# print/LPRng-core/PLIST
@comment $NetBSD: PLIST,v 1.7 2005/05/02 20:34:04 reed Exp $
bin/cancel
bin/lp
bin/lpq
bin/lpr
bin/lprm
bin/lpstat
lib/liblpr.la
libexec/LPRng/lpbanner
libexec/LPRng/lpf
libexec/LPRng/pclbanner
libexec/LPRng/psbanner
man/man1/cancel.1
man/man1/lp.1
man/man1/lpbanner.1
man/man1/lpf.1
man/man1/lpq.1
man/man1/lpr.1
man/man1/lprm.1
man/man1/lprng_certs.1
man/man1/lprng_index_certs.1
man/man1/lpstat.1
man/man1/monitor.1
man/man1/pclbanner.1
man/man1/psbanner.1
man/man5/lpd.conf.5
man/man5/lpd.perms.5
man/man5/printcap.5
man/man8/checkpc.8
man/man8/lpc.8
man/man8/lpd.8
sbin/checkpc
sbin/lpc
sbin/lpd
sbin/lprng_certs
sbin/lprng_index_certs
share/examples/LPRng/lpd.conf.example
share/examples/LPRng/lpd.perms.example
share/examples/LPRng/printcap.example
share/examples/rc.d/LPRng
@dirrm share/examples/LPRng
@dirrm libexec/LPRng
