# net/bind912/files/smf/manifest.xml
<?xml version="1.0"?>
<!--
CDDL HEADER START

The contents of this file are subject to the terms of the
Common Development and Distribution License (the "License").
You may not use this file except in compliance with the License.

You can obtain a copy of the license at usr/src/OPENSOLARIS.LICENSE
or http://www.opensolaris.org/os/licensing.
See the License for the specific language governing permissions
and limitations under the License.

When distributing Covered Code, include this CDDL HEADER in each
file and include the License file at usr/src/OPENSOLARIS.LICENSE.
If applicable, add the following below this CDDL HEADER, with the
fields enclosed by brackets "[]" replaced with your own identifying
information: Portions Copyright [yyyy] [name of copyright owner]

CDDL HEADER END
-->
<!DOCTYPE service_bundle SYSTEM "/usr/share/lib/xml/dtd/service_bundle.dtd.1">
<!--
        Copyright 2007 Sun Microsystems, Inc.  All rights reserved.
        Use is subject to license terms.

        ident   "@(#)server.xml 1.8     07/06/19 SMI"
-->
<!-- Portions Copyright 2007 Vedran Bender -->
<service_bundle type='manifest' name='export'>
  <service name='@SMF_PREFIX@/@SMF_NAME@' type='service' version='1'>
    <dependency name='filesystem_minimal' grouping='require_all' restart_on='none' type='service'>
      <service_fmri value='svc:/system/filesystem/local' />
    </dependency>
    <dependency name='loopback' grouping='require_any' restart_on='error' type='service'>
      <service_fmri value='svc:/network/loopback' />
    </dependency>
    <dependency name='network' grouping='optional_all' restart_on='error' type='service'>
      <service_fmri value='svc:/milestone/network' />
    </dependency>
    <dependency name='config-files' grouping='require_any' restart_on='refresh' type='path'>
      <service_fmri value='file://localhost@PKG_SYSCONFDIR@/named.conf' />
    </dependency>
    <!--
        In order to run multiple named(1M) processes with their own
        configuration file or properties each must have a unique
        instance.
    -->
    <instance name='default' enabled='false'>
      <exec_method type='method' name='start' exec='@PREFIX@/@SMF_METHOD_FILE.named@ %m %i' timeout_seconds='60'>
        <method_context>
          <method_credential user='root' group='root' />
        </method_context>
      </exec_method>
      <exec_method type='method' name='stop' exec='@PREFIX@/@SMF_METHOD_FILE.named@ %m %i %{restarter/contract}' timeout_seconds='60'>
        <method_context>
          <method_credential user='root' group='root' />
        </method_context>

      </exec_method>
      <!--
              SIGHUP causes named to reread its configuration file, but not any
              of the properties below.
      -->
      <exec_method type='method' name='refresh' exec=':kill -HUP' timeout_seconds='60'>
        <method_context />
      </exec_method>
      <property_group name='general' type='framework'>
        <!-- manage DNS server state -->
        <propval name='action_authorization' type='astring' value='solaris.smf.manage.bind' />
        <propval name='value_authorization' type='astring' value='solaris.smf.manage.bind' />
      </property_group>
      <!-- Default property settings for named(1M) instance. -->
      <property_group name='options' type='application'>
        <!--
                server: specifies an alternative server command.  If
                not specified the default /usr/sbin/named is used.
        -->
        <propval name='server' type='astring' value='' />
        <!--
                configuration_file: specifies an alternative
                configuration file to be used. The property is similar
                to named(1M) command line option '-c'
        -->
        <propval name='configuration_file' type='astring' value='' />
        <!--
                ip_interfaces: specifies which IP transport BIND will
                transmit on. Possible values are 'IPv4' or 'IPv6'. Any
                other setting assumes 'all', the default.
                Equivalent command line option '-4' or '-6'.
        -->
        <propval name='ip_interfaces' type='astring' value='all' />
        <!--
                listen_on_port: Specifies the default UDP and TCP port
                which will be used to listen for DNS requests.
                Equivalent command line option '-p <integer>'.
        -->
        <propval name='listen_on_port' type='integer' value='0' />
        <!--
                debug_level: Specifies the default debug level.  The
                default is 0; no debugging. The Higher the number the
                more verbose debug information becomes.
                Equivalent command line option '-d <integer>'.
        -->
        <propval name='debug_level' type='integer' value='0' />
        <!--
                threads: Specifies the number of cpu worker threads to
                create.  The default of 0 causes named to try and
                determine the number of CPUs present and create one
                thread per CPU.
                Equivalent command line option '-n <integer>'.
        -->
        <propval name='threads' type='integer' value='0' />
        <!--
                chroot_dir: Change the root directory using chroot(2)
                to directory after processing the command line
                arguments, but before reading the configuration file.
                Equivalent command line option '-t <pathname>'.
        -->
        <propval name='chroot_dir' type='astring' value='' />
	<!--
		user: Change the user id after processing command line
		arguments, but before reading the configuration file.
		Equivalent command line option '-u <user>'.
	-->
        <propval name='user' type='astring' value='named' />
      </property_group>
    </instance>
    <template>
      <common_name>
        <loctext xml:lang='C'>BIND DNS server</loctext>
      </common_name>
      <documentation>
        <manpage title='named' section='1M' manpath='@PREFIX@/@PKGMANDIR@' />
      </documentation>
    </template>
  </service>
</service_bundle>
