# mail/clamsmtp/PLIST
@comment $NetBSD: PLIST,v 1.5 2005/05/02 20:33:59 reed Exp $
man/man5/clamsmtpd.conf.5
man/man8/clamsmtpd.8
sbin/clamsmtpd
share/examples/clamsmtp/clamsmtpd.conf
share/examples/clamsmtp/virus_action.sh
share/examples/rc.d/clamsmtpd
@dirrm share/examples/clamsmtp
