# mail/clamsmtp/PLIST
@comment $NetBSD: PLIST,v 1.6 2009/06/14 18:04:27 joerg Exp $
man/man5/clamsmtpd.conf.5
man/man8/clamsmtpd.8
sbin/clamsmtpd
share/examples/clamsmtp/clamsmtpd.conf
share/examples/clamsmtp/virus_action.sh
share/examples/rc.d/clamsmtpd
