# net/oidentd/PLIST
@comment $NetBSD: PLIST,v 1.1.1.1 2004/01/03 23:21:44 tron Exp $
etc/rc.d/oidentd
man/man5/oidentd.conf.5
man/man5/oidentd_masq.conf.5
man/man8/oidentd.8
sbin/oidentd
