# mail/milter-greylist/PLIST
@comment $NetBSD: PLIST,v 1.7 2009/06/14 18:04:38 joerg Exp $
bin/milter-greylist
man/man8/milter-greylist.8
man/man5/greylist.conf.5
share/doc/milter-greylist/README
share/examples/milter-greylist/greylist.conf
share/examples/rc.d/milter-greylist
