# security/xdm-krb4/PLIST
@comment $NetBSD: PLIST,v 1.3 2003/03/03 02:53:07 jschauma Exp $
bin/xdm-krb4
etc/rc.d/xdm_krb4
${IMAKE_MAN_DIR}/xdm-krb4.${IMAKE_MANNEWSUFFIX}
share/examples/xdm-krb4/GiveConsole
share/examples/xdm-krb4/TakeConsole
share/examples/xdm-krb4/Xaccess
share/examples/xdm-krb4/Xresources
share/examples/xdm-krb4/Xservers
share/examples/xdm-krb4/Xsession
share/examples/xdm-krb4/Xsetup_0
share/examples/xdm-krb4/xdm-config
@dirrm share/examples/xdm-krb4
