# x11/xorg-server/PLIST.OpenBSD
@comment $NetBSD: PLIST.OpenBSD,v 1.1 2005/01/25 17:55:45 xtraeme Exp $
${X11ROOT_PREFIX}/bin/kbd_mode
${X11ROOT_PREFIX}/etc/init.d/xprint
${X11ROOT_PREFIX}/${MAN_DIR}/kbd_mode.${MAN_SUFFIX}
@dirrm ${X11ROOT_PREFIX}/etc/init.d
