# x11/xorg-server/PLIST.Linux
@comment $NetBSD: PLIST.Linux,v 1.3 2004/10/08 03:35:25 xtraeme Exp $
${X11ROOT_PREFIX}/etc/init.d/xprint
${X11ROOT_PREFIX}/lib/X11/xserver/C/print/models/CANONBJ10E-GS/fonts
${X11ROOT_PREFIX}/lib/X11/xserver/C/print/models/CANONBJ10E-GS/model-config
${X11ROOT_PREFIX}/lib/X11/xserver/C/print/models/PS2PDFspooldir-GS/fonts
${X11ROOT_PREFIX}/lib/X11/xserver/C/print/models/PS2PDFspooldir-GS/model-config
${X11ROOT_PREFIX}/lib/X11/xserver/C/print/models/PS2PDFspooldir-GS/ps2pdf_spooltodir.sh
${X11ROOT_PREFIX}/lib/X11/xserver/C/print/models/GSdefault/fonts
${X11ROOT_PREFIX}/lib/X11/xserver/C/print/models/GSdefault/model-config
${X11ROOT_PREFIX}/lib/modules/drivers/fbdev_drv.o
${X11ROOT_PREFIX}/lib/modules/extensions/libdri.a
${X11ROOT_PREFIX}/lib/modules/drivers/linux/v4l_drv.o
${X11ROOT_PREFIX}/lib/modules/input/aiptek_drv.o
${X11ROOT_PREFIX}/lib/modules/input/linux/ur98_drv.o
${X11ROOT_PREFIX}/lib/modules/linux/libdrm.a
${X11ROOT_PREFIX}/lib/modules/linux/libfbdevhw.a
${X11ROOT_PREFIX}/lib/modules/linux/libint10.a
${X11ROOT_PREFIX}/${KERNMAN_DIR}/aiptek.${KERNMAN_SUFFIX}
${X11ROOT_PREFIX}/${KERNMAN_DIR}/ur98.${KERNMAN_SUFFIX}
${X11ROOT_PREFIX}/${KERNMAN_DIR}/v4l.${KERNMAN_SUFFIX}
@dirrm ${X11ROOT_PREFIX}/lib/X11/xserver/C/print/models/CANONBJ10E-GS
@dirrm ${X11ROOT_PREFIX}/lib/X11/xserver/C/print/models/GSdefault
@dirrm ${X11ROOT_PREFIX}/lib/X11/xserver/C/print/models/PS2PDFspooldir-GS
@dirrm ${X11ROOT_PREFIX}/lib/modules/drivers/linux
@dirrm ${X11ROOT_PREFIX}/lib/modules/input/linux
@dirrm ${X11ROOT_PREFIX}/lib/modules/linux
@dirrm ${X11ROOT_PREFIX}/etc/init.d
