# chat/dircproxy-devel/PLIST
@comment $NetBSD: PLIST,v 1.3 2006/02/19 15:47:35 adrianp Exp $
bin/dircproxy
bin/dircproxy-crypt
man/man1/dircproxy-crypt.1
man/man1/dircproxy.1
share/dircproxy/cronchk.sh
share/dircproxy/log.pl
share/dircproxy/privmsg-log.pl
share/examples/dircproxy/dircproxyrc.sample
share/examples/rc.d/dircproxy
@dirrm share/examples/dircproxy
@dirrm share/dircproxy
