# chat/dircproxy-devel/PLIST
@comment $NetBSD: PLIST,v 1.2 2005/05/02 20:33:58 reed Exp $
bin/dircproxy
bin/dircproxy-crypt
man/man1/dircproxy-crypt.1
man/man1/dircproxy.1
share/dircproxy/cronchk.sh
share/dircproxy/log.pl
share/dircproxy/privmsg-log.pl
share/examples/dircproxy/dircproxyrc
share/examples/rc.d/dircproxy
@dirrm share/examples/dircproxy
@dirrm share/dircproxy
