# multimedia/ccxstream/PLIST
@comment $NetBSD: PLIST,v 1.1.1.1 2009/11/27 08:38:49 wiz Exp $
bin/ccxstream
bin/ccxtest
share/doc/ccxstream/COPYRIGHT
share/doc/ccxstream/ChangeLog
share/doc/ccxstream/PORTING
share/doc/ccxstream/README
share/doc/ccxstream/TODO
share/doc/ccxstream/xbmsp-rc
share/doc/ccxstream/xbmsp-xml.txt
share/doc/ccxstream/xbmsp.txt
share/examples/rc.d/ccxstream
