# www/squid/pkg/PLIST
@comment $NetBSD: PLIST,v 1.1 1998/05/23 21:27:22 tron Exp $
sbin/RunAccel
sbin/RunCache
sbin/cachemgr.cgi
sbin/client
sbin/dnsserver
sbin/ftpget
sbin/pinger
sbin/squid
sbin/unlinkd
etc/rc.d/squid.sh
etc/squid/squid.conf.default
@unexec rm -f %D/etc/squid/squid.conf
@dirrm etc/squid
