# www/squid/pkg/PLIST
@comment $NetBSD: PLIST,v 1.4 2001/01/26 04:56:14 hubertf Exp $
etc/rc.d/squid
libexec/cachemgr.cgi
libexec/dnsserver
libexec/pinger
libexec/unlinkd
sbin/RunAccel
sbin/RunCache
sbin/client
sbin/squid
