# www/squid/pkg/PLIST
@comment $NetBSD: PLIST,v 1.7 2001/09/18 16:24:58 tron Exp $
etc/rc.d/squid
etc/squid/mime.conf.default
etc/squid/squid.conf.default
@unexec rmdir %D/etc/squid 2>/dev/null || true
libexec/cachemgr.cgi
sbin/RunAccel
sbin/RunCache
sbin/client
sbin/squid
