# www/squid/pkg/PLIST
@comment $NetBSD: PLIST,v 1.6 2001/07/29 16:41:19 tron Exp $
etc/rc.d/squid
etc/squid/mime.conf.default
etc/squid/squid.conf.default
@unexec rmdir %D/etc/squid 2>/dev/null || true
libexec/cachemgr.cgi
libexec/pinger
libexec/unlinkd
sbin/RunAccel
sbin/RunCache
sbin/client
sbin/squid
