# www/squid/pkg/PLIST
@comment $NetBSD: PLIST,v 1.3 1998/12/06 12:42:01 tron Exp $
etc/rc.d/squid.sh
libexec/cachemgr.cgi
libexec/dnsserver
libexec/pinger
libexec/unlinkd
sbin/RunAccel
sbin/RunCache
sbin/client
sbin/squid
