# security/prelude-correlator/PLIST
@comment $NetBSD: PLIST,v 1.3 2009/08/24 12:30:09 hasso Exp $
bin/prelude-correlator
sbin/run-prelude-correlator
${PYSITELIB}/PreludeCorrelator/__init__.pyc
${PYSITELIB}/PreludeCorrelator/__init__.pyo
${PYSITELIB}/PreludeCorrelator/config.pyc
${PYSITELIB}/PreludeCorrelator/config.pyo
${PYSITELIB}/PreludeCorrelator/context.pyc
${PYSITELIB}/PreludeCorrelator/context.pyo
${PYSITELIB}/PreludeCorrelator/idmef.pyc
${PYSITELIB}/PreludeCorrelator/idmef.pyo
${PYSITELIB}/PreludeCorrelator/log.pyc
${PYSITELIB}/PreludeCorrelator/log.pyo
${PYSITELIB}/PreludeCorrelator/main.pyc
${PYSITELIB}/PreludeCorrelator/main.pyo
${PYSITELIB}/PreludeCorrelator/pluginmanager.pyc
${PYSITELIB}/PreludeCorrelator/pluginmanager.pyo
${PYSITELIB}/PreludeCorrelator/plugins/__init__.pyc
${PYSITELIB}/PreludeCorrelator/plugins/__init__.pyo
${PYSITELIB}/PreludeCorrelator/plugins/bruteforce.pyc
${PYSITELIB}/PreludeCorrelator/plugins/bruteforce.pyo
${PYSITELIB}/PreludeCorrelator/plugins/businesshour.pyc
${PYSITELIB}/PreludeCorrelator/plugins/businesshour.pyo
${PYSITELIB}/PreludeCorrelator/plugins/dshield.pyc
${PYSITELIB}/PreludeCorrelator/plugins/dshield.pyo
${PYSITELIB}/PreludeCorrelator/plugins/firewall.pyc
${PYSITELIB}/PreludeCorrelator/plugins/firewall.pyo
${PYSITELIB}/PreludeCorrelator/plugins/opensshauth.pyc
${PYSITELIB}/PreludeCorrelator/plugins/opensshauth.pyo
${PYSITELIB}/PreludeCorrelator/plugins/scan.pyc
${PYSITELIB}/PreludeCorrelator/plugins/scan.pyo
${PYSITELIB}/PreludeCorrelator/plugins/worm.pyc
${PYSITELIB}/PreludeCorrelator/plugins/worm.pyo
${PYSITELIB}/PreludeCorrelator/require.pyc
${PYSITELIB}/PreludeCorrelator/require.pyo
${PYSITELIB}/PreludeCorrelator/siteconfig.py
${PYSITELIB}/PreludeCorrelator/siteconfig.pyc
${PYSITELIB}/PreludeCorrelator/siteconfig.pyo
${PYSITELIB}/PreludeCorrelator/utils.pyc
${PYSITELIB}/PreludeCorrelator/utils.pyo
${PYSITELIB}/prelude_correlator-0.9.0_beta6-py2.5.egg-info/PKG-INFO
${PYSITELIB}/prelude_correlator-0.9.0_beta6-py2.5.egg-info/SOURCES.txt
${PYSITELIB}/prelude_correlator-0.9.0_beta6-py2.5.egg-info/dependency_links.txt
${PYSITELIB}/prelude_correlator-0.9.0_beta6-py2.5.egg-info/entry_points.txt
${PYSITELIB}/prelude_correlator-0.9.0_beta6-py2.5.egg-info/top_level.txt
share/examples/rc.d/preludecorrelator
