# security/prelude-correlator/PLIST
@comment $NetBSD: PLIST,v 1.6 2010/03/23 18:01:47 joerg Exp $
bin/prelude-correlator
${PYSITELIB}/PreludeCorrelator/__init__.py
${PYSITELIB}/PreludeCorrelator/__init__.pyc
${PYSITELIB}/PreludeCorrelator/__init__.pyo
${PYSITELIB}/PreludeCorrelator/config.py
${PYSITELIB}/PreludeCorrelator/config.pyc
${PYSITELIB}/PreludeCorrelator/config.pyo
${PYSITELIB}/PreludeCorrelator/context.py
${PYSITELIB}/PreludeCorrelator/context.pyc
${PYSITELIB}/PreludeCorrelator/context.pyo
${PYSITELIB}/PreludeCorrelator/idmef.py
${PYSITELIB}/PreludeCorrelator/idmef.pyc
${PYSITELIB}/PreludeCorrelator/idmef.pyo
${PYSITELIB}/PreludeCorrelator/log.py
${PYSITELIB}/PreludeCorrelator/log.pyc
${PYSITELIB}/PreludeCorrelator/log.pyo
${PYSITELIB}/PreludeCorrelator/main.py
${PYSITELIB}/PreludeCorrelator/main.pyc
${PYSITELIB}/PreludeCorrelator/main.pyo
${PYSITELIB}/PreludeCorrelator/pluginmanager.py
${PYSITELIB}/PreludeCorrelator/pluginmanager.pyc
${PYSITELIB}/PreludeCorrelator/pluginmanager.pyo
${PYSITELIB}/PreludeCorrelator/plugins/__init__.py
${PYSITELIB}/PreludeCorrelator/plugins/__init__.pyc
${PYSITELIB}/PreludeCorrelator/plugins/__init__.pyo
${PYSITELIB}/PreludeCorrelator/plugins/bruteforce.py
${PYSITELIB}/PreludeCorrelator/plugins/bruteforce.pyc
${PYSITELIB}/PreludeCorrelator/plugins/bruteforce.pyo
${PYSITELIB}/PreludeCorrelator/plugins/businesshour.py
${PYSITELIB}/PreludeCorrelator/plugins/businesshour.pyc
${PYSITELIB}/PreludeCorrelator/plugins/businesshour.pyo
${PYSITELIB}/PreludeCorrelator/plugins/dshield.py
${PYSITELIB}/PreludeCorrelator/plugins/dshield.pyc
${PYSITELIB}/PreludeCorrelator/plugins/dshield.pyo
${PYSITELIB}/PreludeCorrelator/plugins/firewall.py
${PYSITELIB}/PreludeCorrelator/plugins/firewall.pyc
${PYSITELIB}/PreludeCorrelator/plugins/firewall.pyo
${PYSITELIB}/PreludeCorrelator/plugins/opensshauth.py
${PYSITELIB}/PreludeCorrelator/plugins/opensshauth.pyc
${PYSITELIB}/PreludeCorrelator/plugins/opensshauth.pyo
${PYSITELIB}/PreludeCorrelator/plugins/scan.py
${PYSITELIB}/PreludeCorrelator/plugins/scan.pyc
${PYSITELIB}/PreludeCorrelator/plugins/scan.pyo
${PYSITELIB}/PreludeCorrelator/plugins/spamhausdrop.py
${PYSITELIB}/PreludeCorrelator/plugins/spamhausdrop.pyc
${PYSITELIB}/PreludeCorrelator/plugins/spamhausdrop.pyo
${PYSITELIB}/PreludeCorrelator/plugins/worm.py
${PYSITELIB}/PreludeCorrelator/plugins/worm.pyc
${PYSITELIB}/PreludeCorrelator/plugins/worm.pyo
${PYSITELIB}/PreludeCorrelator/require.py
${PYSITELIB}/PreludeCorrelator/require.pyc
${PYSITELIB}/PreludeCorrelator/require.pyo
${PYSITELIB}/PreludeCorrelator/siteconfig.py
${PYSITELIB}/PreludeCorrelator/siteconfig.pyc
${PYSITELIB}/PreludeCorrelator/siteconfig.pyo
${PYSITELIB}/PreludeCorrelator/utils.py
${PYSITELIB}/PreludeCorrelator/utils.pyc
${PYSITELIB}/PreludeCorrelator/utils.pyo
${PYSITELIB}/${EGG_INFODIR}/PKG-INFO
${PYSITELIB}/${EGG_INFODIR}/SOURCES.txt
${PYSITELIB}/${EGG_INFODIR}/dependency_links.txt
${PYSITELIB}/${EGG_INFODIR}/entry_points.txt
${PYSITELIB}/${EGG_INFODIR}/not-zip-safe
${PYSITELIB}/${EGG_INFODIR}/top_level.txt
sbin/run-prelude-correlator
share/examples/dshield.dat
share/examples/prelude-correlator.conf
share/examples/rc.d/preludecorrelator
share/examples/spamhaus_drop.dat
