# security/prelude-correlator/PLIST
@comment $NetBSD: PLIST,v 1.1.1.1 2008/07/21 12:16:46 shannonjr Exp $
bin/prelude-correlator
sbin/run-prelude-correlator
include/prelude-correlator/prelude-correlator.h
lib/prelude-correlator/lua.la
share/examples/prelude-correlator/prelude-correlator.conf
share/examples/rc.d/preludecorrelator
share/examples/prelude-correlator/lua-rules/brute-force.lua
share/examples/prelude-correlator/lua-rules/business-hour.lua
share/examples/prelude-correlator/lua-rules/firewall.lua
share/examples/prelude-correlator/lua-rules/scan.lua
share/examples/prelude-correlator/lua-rules/worm.lua
share/prelude-correlator/lua/lib.lua
@dirrm share/prelude-correlator/lua
@dirrm share/examples/prelude-correlator
@dirrm lib/prelude-correlator
@dirrm include/prelude-correlator
