# net/dhcpcd/PLIST
@comment $NetBSD: PLIST,v 1.3 2009/06/14 18:09:25 joerg Exp $
libexec/dhcpcd-hooks/01-test
libexec/dhcpcd-hooks/10-mtu
libexec/dhcpcd-hooks/20-resolv.conf
libexec/dhcpcd-hooks/29-lookup-hostname
libexec/dhcpcd-hooks/30-hostname
libexec/dhcpcd-hooks/50-ntp.conf
libexec/dhcpcd-run-hooks
man/man5/dhcpcd.conf.5
man/man8/dhcpcd-run-hooks.8
man/man8/dhcpcd.8
sbin/dhcpcd
share/examples/dhcpcd.conf
share/examples/rc.d/dhcpcd
