# sysutils/bkpupsd/PLIST
@comment $NetBSD: PLIST,v 1.5 2005/05/02 20:34:05 reed Exp $
libexec/bkpupsd/battery.shutdown
libexec/bkpupsd/line.lost
libexec/bkpupsd/line.recover
libexec/bkpupsd/line.shutdown
share/examples/rc.d/bkpupsd
@dirrm libexec/bkpupsd
man/cat5/upstab.0
man/cat8/bkpupsd.0
man/man5/upstab.5
man/man8/bkpupsd.8
sbin/bkpupsd
share/doc/bkpupsd/README
share/doc/bkpupsd/README.ja
share/doc/bkpupsd/rc.conf.d.sample
share/examples/bkpupsd/upstab
@dirrm share/examples/bkpupsd
@dirrm share/doc/bkpupsd
