# mail/courier-imap/files/smf/manifest.xml
<?xml version='1.0'?>
<!DOCTYPE service_bundle SYSTEM '/usr/share/lib/xml/dtd/service_bundle.dtd.1'>
<service_bundle type='manifest' name='export'>
  <service name='@SMF_PREFIX@/@SMF_NAME@' type='service' version='1'>
    <create_default_instance enabled="false" />
    <single_instance />
    <dependency name='fs-local' grouping='require_all' restart_on='none' type='service'>
      <service_fmri value='svc:/system/filesystem/local' />
    </dependency>
    <dependency name='network-service' grouping='require_all' restart_on='none' type='service'>
      <service_fmri value='svc:/network/service' />
    </dependency>
    <dependency name='system-log' grouping='optional_all' restart_on='none' type='service'>
      <service_fmri value='svc:/system/system-log' />
    </dependency>
    <dependency name='config-files' grouping='require_any' restart_on='refresh' type='path'>
      <service_fmri value='file://localhost@PKG_SYSCONFDIR@/pop3d-ssl' />
      <service_fmri value='file://localhost@PKG_SYSCONFDIR@/pop3d' />
      <service_fmri value='file://localhost@PKG_SYSCONFDIR@/imapd-ssl' />
      <service_fmri value='file://localhost@PKG_SYSCONFDIR@/imapd' />
    </dependency>
    <dependency name='authlib' grouping='require_all' restart_on='none' type='service'>
      <service_fmri value='svc:/@SMF_PREFIX@/courier-authlib' />
    </dependency>
    <exec_method name='start' type='method' exec='@PREFIX@/@SMF_METHOD_FILE.courier-imap@ start' timeout_seconds='120' />
    <exec_method name='stop' type='method' exec='@PREFIX@/@SMF_METHOD_FILE.courier-imap@ stop' timeout_seconds='60' />
    <property_group name='startd' type='framework'>
      <propval name='ignore_error' type='astring' value='core,signal' />
    </property_group>
    <property_group name='general' type='framework'>
      <propval name='action_authorization' type='astring' value='solaris.smf.manage.courier-imap' />
    </property_group>
    <template>
      <common_name>
        <loctext xml:lang='C'>imap - Courier IMAP provides access to Maildirs</loctext>
      </common_name>
      <documentation>
        <manpage title='imapd' section='8' manpath='@PREFIX@/@PKGMANDIR@' />
      </documentation>
    </template>
  </service>
</service_bundle>
