# net/gated/files/gated
#!/bin/sh
#
# $NetBSD: gated,v 1.3 2007/04/24 12:06:18 tron Exp $
#

# PROVIDE: gated
# REQUIRE: network mountcritlocal
# BEFORE:  NETWORKING

. /etc/rc.subr

name="gated"
rcvar=$name
command="@PREFIX@/sbin/${name}"
pidfile="/var/run/${name}.pid"
required_files="/etc/${name}.conf"

load_rc_config $name
run_rc_command "$1"
