# net/gated/files/gated
#!/bin/sh
#
# $NetBSD: gated,v 1.2 2001/04/05 17:01:11 wiz Exp $
#

# PROVIDE: gated
# REQUIRE: DAEMON

. /etc/rc.subr

name="gated"
rcvar=$name
command="@PREFIX@/sbin/${name}"
pidfile="/var/run/${name}.pid"
required_files="/etc/${name}.conf"

load_rc_config $name
run_rc_command "$1"
