# net/ra-rtsp-proxy/pkg/PLIST
@comment $NetBSD: PLIST,v 1.1.1.1 1999/09/01 20:25:03 tron Exp $
etc/rc.d/rtspd.sh
man/man8/rtspd.8
sbin/rtspd
share/doc/rtspd/README
share/doc/rtspd/license.txt
@dirrm share/doc/rtspd
