# net/ra-rtsp-proxy/pkg/PLIST
@comment $NetBSD: PLIST,v 1.2 1999/09/01 21:00:11 tron Exp $
etc/rc.d/rtspd.sh
man/man8/rtspd.8
sbin/rtspd
share/doc/ra-rtsp-proxy/README
share/doc/ra-rtsp-proxy/license.txt
@dirrm share/doc/rtspd
