# sysutils/bacula/PLIST.server
@comment $NetBSD: PLIST.server,v 1.10 2013/01/04 23:25:49 dholland Exp $
lib/libbaccats.la
lib/libbaccats-${BACULA_DB}.la
lib/libbacsql.la
libexec/bacula/chio-changer
libexec/bacula/create_${BACULA_DB}_database
libexec/bacula/create_bacula_database
libexec/bacula/delete_catalog_backup
libexec/bacula/drop_${BACULA_DB}_database
libexec/bacula/drop_${BACULA_DB}_tables
libexec/bacula/drop_bacula_database
libexec/bacula/drop_bacula_tables
libexec/bacula/grant_${BACULA_DB}_privileges
libexec/bacula/grant_bacula_privileges
libexec/bacula/make_${BACULA_DB}_tables
libexec/bacula/make_bacula_tables
libexec/bacula/make_catalog_backup
libexec/bacula/make_catalog_backup.pl
libexec/bacula/query.sql
libexec/bacula/update_${BACULA_DB}_tables
libexec/bacula/update_bacula_tables
sbin/bacula-dir
sbin/bacula-sd
sbin/bcopy
sbin/bextract
sbin/bls
sbin/bregex
sbin/bscan
sbin/bsmtp
sbin/btape
sbin/bwild
sbin/dbcheck
${STATIC}sbin/static-bacula-dir
${STATIC}sbin/static-bacula-sd
share/examples/bacula/bacula-dir.conf
share/examples/bacula/bacula-sd.conf
share/examples/rc.d/bacula
share/examples/rc.d/bacula-dir
share/examples/rc.d/bacula-sd
