# sysutils/bacula/PLIST.common
@comment $NetBSD: PLIST.common,v 1.1 2005/12/08 01:04:44 wiz Exp $
libexec/bacula/bacula
libexec/bacula/bconsole
libexec/bacula/btraceback.dbx
libexec/bacula/btraceback.gdb
libexec/bacula/dvd-handler
libexec/bacula/gconsole
libexec/bacula/mtx-changer
libexec/bacula/startmysql
libexec/bacula/stopmysql
man/man8/bacula.8
sbin/bacula-fd
sbin/bconsole
sbin/btraceback
share/examples/bacula/bacula-fd.conf
share/examples/bacula/bconsole.conf
share/examples/rc.d/bacula-fd
