# sysutils/bacula/PLIST
@comment $NetBSD: PLIST,v 1.2 2005/05/02 20:34:05 reed Exp $
share/examples/bacula/bacula-dir.conf
share/examples/bacula/bacula-fd.conf
share/examples/bacula/bacula-sd.conf
share/examples/bacula/bconsole.conf
libexec/bacula/query.sql
libexec/bacula/bacula
libexec/bacula/bconsole
libexec/bacula/btraceback.gdb
libexec/bacula/btraceback.dbx
libexec/bacula/create_${BACULA_DB}_database
libexec/bacula/create_bacula_database
libexec/bacula/delete_catalog_backup
libexec/bacula/drop_${BACULA_DB}_database
libexec/bacula/drop_${BACULA_DB}_tables
libexec/bacula/drop_bacula_database
libexec/bacula/drop_bacula_tables
libexec/bacula/gconsole
libexec/bacula/grant_${BACULA_DB}_privileges
libexec/bacula/grant_bacula_privileges
libexec/bacula/make_${BACULA_DB}_tables
libexec/bacula/make_bacula_tables
libexec/bacula/make_catalog_backup
libexec/bacula/mtx-changer
libexec/bacula/startmysql
libexec/bacula/stopmysql
libexec/bacula/update_${BACULA_DB}_tables
libexec/bacula/update_bacula_tables
sbin/bacula-dir
sbin/bacula-fd
sbin/bacula-sd
sbin/bconsole
sbin/bcopy
sbin/bextract
sbin/bls
sbin/bscan
sbin/bsmtp
sbin/btape
sbin/btraceback
sbin/dbcheck
share/examples/rc.d/bacula-fd
share/examples/rc.d/bacula-sd
share/examples/rc.d/bacula-dir
share/examples/rc.d/bacula
@dirrm share/examples/bacula
@dirrm libexec/bacula
