# time/anacron/PLIST
@comment $NetBSD: PLIST,v 1.4 2009/06/14 18:17:56 joerg Exp $
man/man5/anacrontab.5
man/man8/anacron.8
sbin/anacron
share/examples/anacron/anacrontab.FreeBSD
share/examples/anacron/anacrontab.NetBSD
share/examples/rc.d/anacron
