# security/stud/PLIST
@comment $NetBSD: PLIST,v 1.1 2013/03/16 19:41:35 jym Exp $
bin/stud
man/man8/stud.8
share/examples/stud/stud.conf
share/examples/rc.d/stud
