# net/argus/PLIST
@comment $NetBSD: PLIST,v 1.2 2005/05/02 20:34:01 reed Exp $
bin/argus-config
lib/argus/Alias.pm
lib/argus/ArgusCtl.pm
lib/argus/ArgusEncode.pm
lib/argus/Artistic
lib/argus/BER.pm
lib/argus/BaseIO.pm
lib/argus/Commands.pm
lib/argus/Conf.pm
lib/argus/Configable.pm
lib/argus/Control.pm
lib/argus/Cron.pm
lib/argus/DNS.pm
lib/argus/DataBase.pm
lib/argus/Doc.pm
lib/argus/Error.pm
lib/argus/Graph.pm
lib/argus/GraphData.pm
lib/argus/Group.pm
lib/argus/Image.pm
lib/argus/LICENSE
lib/argus/LOCATION
lib/argus/MonEl.pm
lib/argus/NotMe.pm
lib/argus/Notify.pm
lib/argus/NullConf.pm
lib/argus/NullCtl.pm
lib/argus/Override.pm
lib/argus/Ping.pm
lib/argus/Prog.pm
lib/argus/README
lib/argus/SNMP.pm
lib/argus/SSL.pm
lib/argus/Self.pm
lib/argus/Server.pm
lib/argus/Service.pm
lib/argus/Stats.pm
lib/argus/TCP.pm
lib/argus/TestPort.pm
lib/argus/UDP.pm
lib/argus/UserCron.pm
lib/argus/Web.pm
lib/argus/conf.pl
lib/argus/graphd
lib/argus/localization.pl
lib/argus/misc.pl
lib/argus/picasso
lib/argus/web_acl.pl
lib/argus/web_auth_file.pl
lib/argus/web_misc.pl
sbin/arguscgi
sbin/argusctl
sbin/argusd
share/examples/rc.d/argusd
@dirrm lib/argus
