# net/argus/PLIST
@comment $NetBSD: PLIST,v 1.3 2007/03/18 18:45:14 adam Exp $
bin/argus-config
lib/argus/Alias.pm
lib/argus/Argus/Agent.pm
lib/argus/Argus/Asterisk.pm
lib/argus/Argus/BER.pm
lib/argus/Argus/Ctl.pm
lib/argus/Argus/Encode.pm
lib/argus/Argus/IP.pm
lib/argus/Argus/RPC.pm
lib/argus/Argus/RPC/TCP.pm
lib/argus/Argus/RPC/UDP.pm
lib/argus/Argus/SIP.pm
lib/argus/Argus/SIP/TCP.pm
lib/argus/Argus/SIP/UDP.pm
lib/argus/Artistic
lib/argus/BaseIO.pm
lib/argus/Chart/Strip.pm
lib/argus/Commands.pm
lib/argus/Conf.pm
lib/argus/Configable.pm
lib/argus/Control.pm
lib/argus/Cron.pm
lib/argus/DNS.pm
lib/argus/DNS/TCP.pm
lib/argus/DNS/UDP.pm
lib/argus/DataBase.pm
lib/argus/Doc.pm
lib/argus/Error.pm
lib/argus/Graph.pm
lib/argus/GraphData.pm
lib/argus/Group.pm
lib/argus/LICENSE
lib/argus/LOCATION
lib/argus/MonEl.pm
lib/argus/NotMe.pm
lib/argus/Notify.pm
lib/argus/NullConf.pm
lib/argus/NullCtl.pm
lib/argus/Override.pm
lib/argus/Ping.pm
lib/argus/Prog.pm
lib/argus/README
lib/argus/Resolv.pm
lib/argus/SNMP.pm
lib/argus/SSL.pm
lib/argus/Self.pm
lib/argus/Server.pm
lib/argus/Service.pm
lib/argus/Stats.pm
lib/argus/TCP.pm
lib/argus/TestPort.pm
lib/argus/UDP.pm
lib/argus/UserCron.pm
lib/argus/Web.pm
lib/argus/conf.pl
lib/argus/graphd
lib/argus/localization.pl
lib/argus/misc.pl
lib/argus/picasso
lib/argus/web_acl.pl
lib/argus/web_auth_file.pl
lib/argus/web_misc.pl
sbin/argus-agent
sbin/arguscgi
sbin/argusctl
sbin/argusd
share/examples/rc.d/argusd
@dirrm lib/argus/DNS
@dirrm lib/argus/Chart
@dirrm lib/argus/Argus/SIP
@dirrm lib/argus/Argus/RPC
@dirrm lib/argus/Argus
@dirrm lib/argus
