# net/argus/PLIST
@comment $NetBSD: PLIST,v 1.5 2009/06/14 18:09:21 joerg Exp $
bin/argus-config
etc/argus/config.example
etc/argus/locale/de_DE
etc/argus/locale/el_gr
etc/argus/locale/eubonics
etc/argus/locale/piglatin
etc/argus/locale/sp_ar
etc/argus/locale/template
etc/argus/log
etc/argus/notno
etc/argus/users.example
lib/argus/Alias.pm
lib/argus/Argus/Agent.pm
lib/argus/Argus/Asterisk.pm
lib/argus/Argus/BER.pm
lib/argus/Argus/Ctl.pm
lib/argus/Argus/Encode.pm
lib/argus/Argus/IP.pm
lib/argus/Argus/RPC.pm
lib/argus/Argus/RPC/TCP.pm
lib/argus/Argus/RPC/UDP.pm
lib/argus/Argus/SIP.pm
lib/argus/Argus/SIP/TCP.pm
lib/argus/Argus/SIP/UDP.pm
lib/argus/Artistic
lib/argus/BaseIO.pm
lib/argus/Chart/Strip.pm
lib/argus/Commands.pm
lib/argus/Conf.pm
lib/argus/Configable.pm
lib/argus/Control.pm
lib/argus/Cron.pm
lib/argus/DNS.pm
lib/argus/DNS/TCP.pm
lib/argus/DNS/UDP.pm
lib/argus/DataBase.pm
lib/argus/Doc.pm
lib/argus/Error.pm
lib/argus/Graph.pm
lib/argus/GraphData.pm
lib/argus/Group.pm
lib/argus/LICENSE
lib/argus/LOCATION
lib/argus/MonEl.pm
lib/argus/NotMe.pm
lib/argus/Notify.pm
lib/argus/NullConf.pm
lib/argus/NullCtl.pm
lib/argus/Override.pm
lib/argus/Ping.pm
lib/argus/Prog.pm
lib/argus/README
lib/argus/Resolv.pm
lib/argus/SNMP.pm
lib/argus/SSL.pm
lib/argus/Self.pm
lib/argus/Server.pm
lib/argus/Service.pm
lib/argus/Stats.pm
lib/argus/TCP.pm
lib/argus/TestPort.pm
lib/argus/UDP.pm
lib/argus/UserCron.pm
lib/argus/Web.pm
lib/argus/conf.pl
lib/argus/graphd
lib/argus/localization.pl
lib/argus/misc.pl
lib/argus/picasso
lib/argus/web_acl.pl
lib/argus/web_auth_file.pl
lib/argus/web_misc.pl
sbin/argus-agent
sbin/arguscgi
sbin/argusctl
sbin/argusd
share/examples/rc.d/argusd
@exec ${MKDIR} %D/etc/argus/stats
@exec ${MKDIR} %D/etc/argus/perl
@exec ${MKDIR} %D/etc/argus/notify
@exec ${MKDIR} %D/etc/argus/html
@exec ${MKDIR} %D/etc/argus/gdata
@exec ${MKDIR} %D/etc/argus/gcache
