# textproc/ndtpd/PLIST
@comment $NetBSD: PLIST,v 1.4 2004/03/11 13:51:25 seb Exp $
etc/rc.d/ndtpd
libexec/ndtpstat
share/examples/ndtpd/ndtpd.conf.default
sbin/ndtpd
sbin/ndtpcheck
sbin/ndtpcontrol
sbin/ndtpdaily
sbin/ndtpupgrade
@dirrm share/examples/ndtpd
