# textproc/ndtpd/PLIST
@comment $NetBSD: PLIST,v 1.7 2006/04/07 19:54:17 jlam Exp $
info/ndtpd-ja.info
info/ndtpd.info
libexec/ndtpstat
share/examples/ndtpd/ndtpd.conf.default
sbin/ndtpd
sbin/ndtpcheck
sbin/ndtpcontrol
sbin/ndtpdaily
sbin/ndtpupgrade
share/examples/rc.d/ndtpd
@dirrm share/examples/ndtpd
