# textproc/ndtpd/PLIST
@comment $NetBSD: PLIST,v 1.2 2003/04/13 01:29:42 salo Exp $
etc/rc.d/ndtpd
@unexec install-info --delete %D/info/ndtpd-ja.info %D/info/dir
@unexec install-info --delete %D/info/ndtpd.info %D/info/dir
info/ndtpd-ja.info
info/ndtpd-ja.info-1
info/ndtpd-ja.info-2
info/ndtpd.info
info/ndtpd.info-1
info/ndtpd.info-2
@exec install-info %D/info/ndtpd-ja.info %D/info/dir
@exec install-info %D/info/ndtpd.info %D/info/dir
libexec/ndtpstat
share/examples/ndtpd/ndtpd.conf.default
sbin/ndtpd
sbin/ndtpcheck
sbin/ndtpcontrol
sbin/ndtpdaily
sbin/ndtpupgrade
@dirrm share/examples/ndtpd
