# textproc/ndtpd/PLIST
@comment $NetBSD: PLIST,v 1.1.1.1 2002/03/15 07:19:58 uebayasi Exp $
etc/ndtpd.conf.sample
etc/rc.d/ndtpd.sh
@unexec install-info --delete %D/info/ndtpd-ja.info %D/info/dir
@unexec install-info --delete %D/info/ndtpd.info %D/info/dir
info/ndtpd-ja.info
info/ndtpd-ja.info-1
info/ndtpd-ja.info-2
info/ndtpd.info
info/ndtpd.info-1
info/ndtpd.info-2
@exec install-info %D/info/ndtpd-ja.info %D/info/dir
@exec install-info %D/info/ndtpd.info %D/info/dir
libexec/ndtpstat
sbin/ndtpd
sbin/ndtpcheck
sbin/ndtpcontrol
sbin/ndtpdaily
sbin/ndtpupgrade
