# security/py-prewikka/PLIST
@comment $NetBSD: PLIST,v 1.1.1.1 2006/01/29 16:00:44 shannonjr Exp $
bin/prewikka-httpd
${PYSITELIB}/__init__.py
${PYSITELIB}/__init__.pyc
${PYSITELIB}/__init__.pyo
${PYSITELIB}/prewikka/Auth.pyc
${PYSITELIB}/prewikka/Auth.pyo
${PYSITELIB}/prewikka/CheetahFilters.pyc
${PYSITELIB}/prewikka/CheetahFilters.pyo
${PYSITELIB}/prewikka/Config.pyc
${PYSITELIB}/prewikka/Config.pyo
${PYSITELIB}/prewikka/Core.pyc
${PYSITELIB}/prewikka/Core.pyo
${PYSITELIB}/prewikka/DataSet.pyc
${PYSITELIB}/prewikka/DataSet.pyo
${PYSITELIB}/prewikka/Database.pyc
${PYSITELIB}/prewikka/Database.pyo
${PYSITELIB}/prewikka/Error.pyc
${PYSITELIB}/prewikka/Error.pyo
${PYSITELIB}/prewikka/Filter.pyc
${PYSITELIB}/prewikka/Filter.pyo
${PYSITELIB}/prewikka/IDMEFDatabase.pyc
${PYSITELIB}/prewikka/IDMEFDatabase.pyo
${PYSITELIB}/prewikka/Log.pyc
${PYSITELIB}/prewikka/Log.pyo
${PYSITELIB}/prewikka/ModPythonHandler.pyc
${PYSITELIB}/prewikka/ModPythonHandler.pyo
${PYSITELIB}/prewikka/MyConfigParser.pyc
${PYSITELIB}/prewikka/MyConfigParser.pyo
${PYSITELIB}/prewikka/ParametersNormalizer.pyc
${PYSITELIB}/prewikka/ParametersNormalizer.pyo
${PYSITELIB}/prewikka/Request.pyc
${PYSITELIB}/prewikka/Request.pyo
${PYSITELIB}/prewikka/User.pyc
${PYSITELIB}/prewikka/User.pyo
${PYSITELIB}/prewikka/__init__.py
${PYSITELIB}/prewikka/__init__.pyc
${PYSITELIB}/prewikka/__init__.pyo
${PYSITELIB}/prewikka/modules/__init__.pyc
${PYSITELIB}/prewikka/modules/__init__.pyo
${PYSITELIB}/prewikka/modules/auth/__init__.pyc
${PYSITELIB}/prewikka/modules/auth/__init__.pyo
${PYSITELIB}/prewikka/modules/auth/loginpassword/__init__.pyc
${PYSITELIB}/prewikka/modules/auth/loginpassword/__init__.pyo
${PYSITELIB}/prewikka/modules/auth/loginpassword/loginpassword.pyc
${PYSITELIB}/prewikka/modules/auth/loginpassword/loginpassword.pyo
${PYSITELIB}/prewikka/modules/log/__init__.pyc
${PYSITELIB}/prewikka/modules/log/__init__.pyo
${PYSITELIB}/prewikka/modules/log/stderr/__init__.pyc
${PYSITELIB}/prewikka/modules/log/stderr/__init__.pyo
${PYSITELIB}/prewikka/modules/log/stderr/stderr.pyc
${PYSITELIB}/prewikka/modules/log/stderr/stderr.pyo
${PYSITELIB}/prewikka/siteconfig.py
${PYSITELIB}/prewikka/siteconfig.pyc
${PYSITELIB}/prewikka/siteconfig.pyo
${PYSITELIB}/prewikka/templates/About.py
${PYSITELIB}/prewikka/templates/About.pyc
${PYSITELIB}/prewikka/templates/About.pyo
${PYSITELIB}/prewikka/templates/AlertListing.py
${PYSITELIB}/prewikka/templates/AlertListing.pyc
${PYSITELIB}/prewikka/templates/AlertListing.pyo
${PYSITELIB}/prewikka/templates/ClassicLayout.py
${PYSITELIB}/prewikka/templates/ClassicLayout.pyc
${PYSITELIB}/prewikka/templates/ClassicLayout.pyo
${PYSITELIB}/prewikka/templates/Command.py
${PYSITELIB}/prewikka/templates/Command.pyc
${PYSITELIB}/prewikka/templates/Command.pyo
${PYSITELIB}/prewikka/templates/ErrorTemplate.py
${PYSITELIB}/prewikka/templates/ErrorTemplate.pyc
${PYSITELIB}/prewikka/templates/ErrorTemplate.pyo
${PYSITELIB}/prewikka/templates/FilterEdition.py
${PYSITELIB}/prewikka/templates/FilterEdition.pyc
${PYSITELIB}/prewikka/templates/FilterEdition.pyo
${PYSITELIB}/prewikka/templates/HTMLDocument.py
${PYSITELIB}/prewikka/templates/HTMLDocument.pyc
${PYSITELIB}/prewikka/templates/HTMLDocument.pyo
${PYSITELIB}/prewikka/templates/HeartbeatAnalyze.py
${PYSITELIB}/prewikka/templates/HeartbeatAnalyze.pyc
${PYSITELIB}/prewikka/templates/HeartbeatAnalyze.pyo
${PYSITELIB}/prewikka/templates/HeartbeatListing.py
${PYSITELIB}/prewikka/templates/HeartbeatListing.pyc
${PYSITELIB}/prewikka/templates/HeartbeatListing.pyo
${PYSITELIB}/prewikka/templates/LoginPasswordForm.py
${PYSITELIB}/prewikka/templates/LoginPasswordForm.pyc
${PYSITELIB}/prewikka/templates/LoginPasswordForm.pyo
${PYSITELIB}/prewikka/templates/MessageDetails.py
${PYSITELIB}/prewikka/templates/MessageDetails.pyc
${PYSITELIB}/prewikka/templates/MessageDetails.pyo
${PYSITELIB}/prewikka/templates/MessageListing.py
${PYSITELIB}/prewikka/templates/MessageListing.pyc
${PYSITELIB}/prewikka/templates/MessageListing.pyo
${PYSITELIB}/prewikka/templates/MessageSummary.py
${PYSITELIB}/prewikka/templates/MessageSummary.pyc
${PYSITELIB}/prewikka/templates/MessageSummary.pyo
${PYSITELIB}/prewikka/templates/PropertiesChangeForm.py
${PYSITELIB}/prewikka/templates/PropertiesChangeForm.pyc
${PYSITELIB}/prewikka/templates/PropertiesChangeForm.pyo
${PYSITELIB}/prewikka/templates/SensorAlertListing.py
${PYSITELIB}/prewikka/templates/SensorAlertListing.pyc
${PYSITELIB}/prewikka/templates/SensorAlertListing.pyo
${PYSITELIB}/prewikka/templates/SensorHeartbeatListing.py
${PYSITELIB}/prewikka/templates/SensorHeartbeatListing.pyc
${PYSITELIB}/prewikka/templates/SensorHeartbeatListing.pyo
${PYSITELIB}/prewikka/templates/SensorListing.py
${PYSITELIB}/prewikka/templates/SensorListing.pyc
${PYSITELIB}/prewikka/templates/SensorListing.pyo
${PYSITELIB}/prewikka/templates/TopLayout.py
${PYSITELIB}/prewikka/templates/TopLayout.pyc
${PYSITELIB}/prewikka/templates/TopLayout.pyo
${PYSITELIB}/prewikka/templates/UserListing.py
${PYSITELIB}/prewikka/templates/UserListing.pyc
${PYSITELIB}/prewikka/templates/UserListing.pyo
${PYSITELIB}/prewikka/templates/UserSettings.py
${PYSITELIB}/prewikka/templates/UserSettings.pyc
${PYSITELIB}/prewikka/templates/UserSettings.pyo
${PYSITELIB}/prewikka/templates/__init__.py
${PYSITELIB}/prewikka/templates/__init__.pyc
${PYSITELIB}/prewikka/templates/__init__.pyo
${PYSITELIB}/prewikka/templates/utils.py
${PYSITELIB}/prewikka/templates/utils.pyc
${PYSITELIB}/prewikka/templates/utils.pyo
${PYSITELIB}/prewikka/utils.pyc
${PYSITELIB}/prewikka/utils.pyo
${PYSITELIB}/prewikka/view.pyc
${PYSITELIB}/prewikka/view.pyo
${PYSITELIB}/prewikka/views/__init__.pyc
${PYSITELIB}/prewikka/views/__init__.pyo
${PYSITELIB}/prewikka/views/commands.pyc
${PYSITELIB}/prewikka/views/commands.pyo
${PYSITELIB}/prewikka/views/filter.pyc
${PYSITELIB}/prewikka/views/filter.pyo
${PYSITELIB}/prewikka/views/messagedetails.pyc
${PYSITELIB}/prewikka/views/messagedetails.pyo
${PYSITELIB}/prewikka/views/messagelisting.pyc
${PYSITELIB}/prewikka/views/messagelisting.pyo
${PYSITELIB}/prewikka/views/messagesummary.pyc
${PYSITELIB}/prewikka/views/messagesummary.pyo
${PYSITELIB}/prewikka/views/misc.pyc
${PYSITELIB}/prewikka/views/misc.pyo
${PYSITELIB}/prewikka/views/sensor.pyc
${PYSITELIB}/prewikka/views/sensor.pyo
${PYSITELIB}/prewikka/views/usermanagement.pyc
${PYSITELIB}/prewikka/views/usermanagement.pyo
sbin/runPrewikka
share/examples/rc.d/prewikka
share/prewikka/database/sqlite.sql
@dirrm share/prewikka/database
@dirrm ${PYSITELIB}/prewikka/views
@dirrm ${PYSITELIB}/prewikka/templates
@dirrm ${PYSITELIB}/prewikka/modules/log/stderr
@dirrm ${PYSITELIB}/prewikka/modules/log
@dirrm ${PYSITELIB}/prewikka/modules/auth/loginpassword
@dirrm ${PYSITELIB}/prewikka/modules/auth
@dirrm ${PYSITELIB}/prewikka/modules
@dirrm ${PYSITELIB}/prewikka
