# mail/mimedefang/PLIST
@comment $NetBSD: PLIST,v 1.6 2005/05/02 20:34:00 reed Exp $
bin/md-mx-ctrl
bin/mimedefang
bin/mimedefang-multiplexor
bin/mimedefang.pl
bin/watch-mimedefang
man/man5/mimedefang-filter.5
man/man7/mimedefang-protocol.7
man/man8/md-mx-ctrl.8
man/man8/mimedefang-multiplexor.8
man/man8/mimedefang.8
man/man8/mimedefang.pl.8
man/man8/watch-mimedefang.8
share/examples/mimedefang/mimedefang-filter
share/examples/mimedefang/sa-mimedefang.cf
share/examples/rc.d/mimedefang
share/examples/rc.d/mimedefang-multiplexor
share/mimedefang/gen-ip-validator.pl
@dirrm share/mimedefang
@dirrm share/examples/mimedefang
