# net/wide-dhcpv6/PLIST
@comment $NetBSD: PLIST,v 1.3 2009/06/14 18:09:48 joerg Exp $
man/man5/dhcp6c.conf.5
man/man5/dhcp6s.conf.5
man/man8/dhcp6c.8
man/man8/dhcp6ctl.8
man/man8/dhcp6relay.8
man/man8/dhcp6s.8
sbin/dhcp6c
sbin/dhcp6ctl
sbin/dhcp6relay
sbin/dhcp6s
share/examples/rc.d/dhcp6c
share/examples/rc.d/dhcp6relay
share/examples/rc.d/dhcp6s
share/examples/wide-dhcpv6/dhcp6c.conf
share/examples/wide-dhcpv6/dhcp6s.conf
