# mail/policyd/PLIST
@comment $NetBSD: PLIST,v 1.1.1.1 2007/08/07 12:34:56 ghen Exp $
libexec/policyd/cleanup
libexec/policyd/policyd
libexec/policyd/stats
share/doc/policyd/README
share/examples/policyd/DATABASE.mysql
share/examples/policyd/blacklist_helo.sql
share/examples/policyd/policyd.conf
share/examples/policyd/whitelist.sql
share/examples/rc.d/policyd
@dirrm share/examples/policyd
@dirrm share/doc/policyd
