# devel/distcc/PLIST
@comment $NetBSD: PLIST,v 1.11 2005/05/02 20:33:59 reed Exp $
bin/distcc
bin/distccd
bin/distccmon-text
man/man1/distcc.1
man/man1/distccd.1
man/man1/distccmon-text.1
share/doc/distcc/AUTHORS
share/doc/distcc/COPYING
share/doc/distcc/INSTALL
share/doc/distcc/NEWS
share/doc/distcc/README
share/doc/distcc/TODO
share/doc/distcc/example/init
share/doc/distcc/example/init-suse
share/doc/distcc/example/logrotate
share/doc/distcc/example/xinetd
share/doc/distcc/protocol-1.txt
share/doc/distcc/protocol-2.txt
share/doc/distcc/reporting-bugs.txt
share/doc/distcc/status-1.txt
share/doc/distcc/survey.txt
share/examples/rc.d/distccd
@dirrm share/doc/distcc/example
@dirrm share/doc/distcc
