# net/iplog/PLIST
@comment $NetBSD: PLIST,v 1.2 2003/07/30 15:41:51 jmmv Exp $
etc/rc.d/iplog
man/man5/iplog.conf.5
man/man8/iplog.8
sbin/iplog
share/examples/iplog/iplog.conf
@dirrm share/examples/iplog
