# net/iplog/PLIST
@comment $NetBSD: PLIST,v 1.5 2009/06/14 18:09:30 joerg Exp $
man/man5/iplog.conf.5
man/man8/iplog.8
sbin/iplog
share/examples/iplog/iplog.conf
share/examples/rc.d/iplog
