# www/nginx/PLIST
@comment $NetBSD: PLIST,v 1.2 2009/06/14 22:00:29 joerg Exp $
sbin/nginx
share/examples/nginx/conf/fastcgi_params
share/examples/nginx/conf/koi-utf
share/examples/nginx/conf/koi-win
share/examples/nginx/conf/mime.types
share/examples/nginx/conf/nginx.conf
share/examples/nginx/conf/win-utf
share/examples/nginx/html/50x.html
share/examples/nginx/html/index.html
share/examples/rc.d/nginx
