# mail/mailfront/PLIST
@comment $NetBSD: PLIST,v 1.1.1.1 2005/07/20 03:30:19 schmonz Exp $
bin/imapfront-auth
bin/pop3front-auth
bin/pop3front-maildir
bin/qmqpfront-echo
bin/qmqpfront-qmail
bin/qmtpfront-echo
bin/qmtpfront-qmail
bin/smtpfront-echo
bin/smtpfront-qmail
bin/smtpfront-reject
share/examples/rc.d/smtpfront
