# net/arpwatch/PLIST
@comment $NetBSD: PLIST,v 1.6 2009/06/14 18:09:21 joerg Exp $
sbin/arpwatch
sbin/arpsnmp
share/arpwatch/ethercodes.dat
share/examples/arpwatch/bihourly.sh
share/examples/arpwatch/arpfetch
man/man8/arpwatch.8
man/man8/arpsnmp.8
share/examples/rc.d/arpwatch
@exec ${MKDIR} ${VARBASE}/db/arpwatch
@exec /usr/bin/touch ${VARBASE}/db/arpwatch/arp.dat
