# net/arpwatch/PLIST
@comment $NetBSD: PLIST,v 1.5 2008/04/04 19:42:15 jwise Exp $
sbin/arpwatch
sbin/arpsnmp
share/arpwatch/ethercodes.dat
share/examples/arpwatch/bihourly.sh
share/examples/arpwatch/arpfetch
man/man8/arpwatch.8
man/man8/arpsnmp.8
share/examples/rc.d/arpwatch
@exec ${MKDIR} ${VARBASE}/db/arpwatch
@exec /usr/bin/touch ${VARBASE}/db/arpwatch/arp.dat
@dirrm share/arpwatch
@dirrm share/examples/arpwatch
