# mail/spamd/PLIST
@comment $NetBSD: PLIST,v 1.3 2009/06/14 18:04:43 joerg Exp $
libexec/spamd
libexec/spamd-setup
libexec/spamlogd
man/cat5/spamd.conf.0
man/cat8/spamd-setup.0
man/cat8/spamd.0
man/cat8/spamdb.0
man/cat8/spamlogd.0
man/man5/spamd.conf.5
man/man8/spamd-setup.8
man/man8/spamd.8
man/man8/spamdb.8
man/man8/spamlogd.8
sbin/spamdb
share/examples/rc.d/pfspamd
share/examples/spamd/spamd.conf
