# sysutils/pluged/PLIST
@comment $NetBSD: PLIST,v 1.3 2014/05/06 15:10:13 jperkin Exp $
sbin/pluged
share/examples/pluged/pluged.conf
${PLIST.rcd}share/examples/rc.d/pluged
