# security/openssh/pkg/PLIST.SunOS
@comment $NetBSD: PLIST.SunOS,v 1.2 2001/02/21 18:55:37 tron Exp $
bin/scp
bin/slogin
bin/ssh
bin/ssh-add
bin/ssh-agent
bin/ssh-keygen
etc/rc.d/sshd
libexec/sftp-server
man/man1/scp.1
man/man1/slogin.1
man/man1/ssh-add.1
man/man1/ssh-agent.1
man/man1/ssh-keygen.1
man/man1/ssh.1
man/man8/sftp-server.8
man/man8/sshd.8
sbin/sshd
share/examples/openssh/ssh_config
share/examples/openssh/ssh_prng_cmds
share/examples/openssh/sshd_config
@dirrm share/examples/ssh
