# security/openssh/pkg/PLIST
@comment $NetBSD: PLIST,v 1.4 2000/07/22 08:22:03 jlam Exp $
bin/scp
bin/slogin
bin/ssh
bin/ssh-add
bin/ssh-agent
bin/ssh-keygen
etc/rc.d/sshd.sh
man/man1/ssh.1
man/man1/scp.1
man/man1/slogin.1
man/man1/ssh-add.1
man/man1/ssh-agent.1
man/man1/ssh-keygen.1
man/man8/sshd.8
sbin/sshd
share/examples/ssh/ssh_config
share/examples/ssh/sshd_config
@dirrm share/examples/ssh
