# security/openssh/pkg/PLIST
@comment $NetBSD: PLIST,v 1.3 2000/05/31 15:32:44 itojun Exp $
bin/ssh
bin/scp
bin/slogin
bin/ssh-add
bin/ssh-agent
bin/ssh-keygen
sbin/sshd
etc/rc.d/sshd.sh
man/man1/ssh.1
man/man1/scp.1
man/man1/slogin.1
man/man1/ssh-add.1
man/man1/ssh-agent.1
man/man1/ssh-keygen.1
man/man8/sshd.8
share/examples/ssh/ssh_config
@comment share/examples/ssh/ssh_prng_cmds
share/examples/ssh/sshd_config
@exec if [ ! -d @SSH_CONF_DIR@ ]; then echo "Creating directory @SSH_CONF_DIR@ for ssh config files.." ; mkdir -p @SSH_CONF_DIR@; fi
@exec if [ ! -f @SSH_CONF_DIR@/ssh_config ]; then echo "Installing example ssh_config in @SSH_CONF_DIR@.." ; /usr/bin/install -c -o root -g wheel -m 0644 %D/share/examples/ssh/ssh_config @SSH_CONF_DIR@; fi
@exec if [ ! -f @SSH_CONF_DIR@/sshd_config ]; then echo "Installing example sshd_config in @SSH_CONF_DIR@.." ; /usr/bin/install -c -o root -g wheel -m 0644 %D/share/examples/ssh/sshd_config @SSH_CONF_DIR@; fi
@comment @exec if [ ! -f @SSH_CONF_DIR@/ssh_prng_cmds ]; then echo "Installing example ssh_prng_cmds in @SSH_CONF_DIR@.." ; /usr/bin/install -c -o root -g wheel -m 0644 %D/share/examples/ssh/ssh_prng_cmds @SSH_CONF_DIR@; fi
@exec if [ ! -f @SSH_CONF_DIR@/ssh_host_key ]; then echo "Generating an RSA secret host key in @SSH_CONF_DIR@.." ; %D/bin/ssh-keygen -N "" -f @SSH_CONF_DIR@/ssh_host_key; fi
@exec if [ ! -f @SSH_CONF_DIR@/ssh_host_dsa_key ]; then echo "Generating a DSA secret host key in @SSH_CONF_DIR@.." ; %D/bin/ssh-keygen -d -N "" -f @SSH_CONF_DIR@/ssh_host_dsa_key; fi
@dirrm share/examples/ssh
