# net/dante/PLIST
@comment $NetBSD: PLIST,v 1.6 2007/12/27 14:24:00 joerg Exp $
bin/socksify
include/dante/socks.h
lib/libdsocks.la
lib/libsocks.la
man/man5/sockd.conf.5
man/man5/socks.conf.5
man/man8/sockd.8
sbin/sockd
share/examples/dante/sockd.conf
share/examples/dante/socks-simple.conf
share/examples/dante/socks.conf
share/examples/rc.d/sockd
@dirrm share/examples/dante
@dirrm include/dante
