# x11/xdm/PLIST
@comment $NetBSD: PLIST,v 1.2 2009/06/14 18:25:11 joerg Exp $
bin/xdm
bin/xdmshell
lib/X11/app-defaults/Chooser
lib/X11/xdm/GiveConsole
lib/X11/xdm/TakeConsole
lib/X11/xdm/Xaccess
lib/X11/xdm/Xreset
lib/X11/xdm/Xresources
lib/X11/xdm/Xservers
lib/X11/xdm/Xsession
lib/X11/xdm/Xsetup_0
lib/X11/xdm/Xstartup
lib/X11/xdm/Xwilling
lib/X11/xdm/chooser
lib/X11/xdm/libXdmGreet.la
lib/X11/xdm/pixmaps/xorg-bw.xpm
lib/X11/xdm/pixmaps/xorg.xpm
lib/X11/xdm/xdm-config
man/man1/xdm.1
share/examples/rc.d/xdm
