# japanese/canna-server/pkg/PLIST
@comment $NetBSD: PLIST,v 1.2 1999/08/27 04:06:59 sakamoto Exp $
@cwd ${DICT_DIR}
@cwd ${DICT_DIR}
canna/canna/bushu.cld
canna/canna/chimei.ctd
canna/canna/dics.dir
canna/canna/fuzokugo.cbd
canna/canna/fuzokugo.cld
canna/canna/hojomwd.ctd
canna/canna/hojoswd.ctd
canna/canna/iroha.cbd
canna/canna/iroha.cld
canna/canna/kanasmpl.ctd
canna/canna/katakana.ctd
canna/canna/keishiki.ctd
canna/canna/necgaiji.ctd
canna/canna/number.ctd
canna/canna/software.ctd
canna/canna/suffix.ctd
@cwd ${PREFIX}
bin/canlisp
bin/cannacheck
bin/cannakill
bin/cannastat
bin/catdic
bin/crfreq
bin/crxdic
bin/crxgram
bin/cshost
bin/ctow
bin/dicar
bin/dpbindic
bin/dpromdic
bin/dpxdic
bin/forcpp
bin/forsort
bin/kpdic
bin/mergeword
bin/mkbindic
bin/mkromdic
bin/splitword
bin/wtoc
etc/rc.d/canna.sh
man/ja_JP.EUC/man1/addwords.1
man/ja_JP.EUC/man1/cannacheck.1
man/ja_JP.EUC/man1/cannakill.1
man/ja_JP.EUC/man1/cannaserver.1
man/ja_JP.EUC/man1/cannastat.1
man/ja_JP.EUC/man1/catdic.1
man/ja_JP.EUC/man1/cpdic.1
man/ja_JP.EUC/man1/cshost.1
man/ja_JP.EUC/man1/ctow.1
man/ja_JP.EUC/man1/delwords.1
man/ja_JP.EUC/man1/dicar.1
man/ja_JP.EUC/man1/dpbindic.1
man/ja_JP.EUC/man1/dpromdic.1
man/ja_JP.EUC/man1/lsdic.1
man/ja_JP.EUC/man1/mkbindic.1
man/ja_JP.EUC/man1/mkdic.1
man/ja_JP.EUC/man1/mkromdic.1
man/ja_JP.EUC/man1/mvdic.1
man/ja_JP.EUC/man1/rmdic.1
man/ja_JP.EUC/man1/splitword.1
man/ja_JP.EUC/man1/wtoc.1
sbin/cannaserver
@exec mkdir ${CANNA_SPOOL}
@exec chown ${CANNAOWNER}:${CANNAGROUP} ${CANNA_SPOOL}
@exec chmod ${CANNA_MODE} ${CANNA_SPOOL}
@unexec rmdir ${DICT_DIR}/canna/canna
@unexec rmdir ${CANNA_SPOOL}
