# chat/ircservices/PLIST
@comment $NetBSD: PLIST,v 1.6 2005/06/21 16:18:52 jlam Exp $
lib/services/convert-db
lib/services/helpfiles/help
lib/services/languages/de
lib/services/languages/en_us
lib/services/languages/es
lib/services/languages/fr
lib/services/languages/hu
lib/services/languages/it
lib/services/languages/ja_euc
lib/services/languages/ja_sjis
lib/services/languages/nl
lib/services/languages/pt
lib/services/languages/ru
lib/services/languages/tr
man/man8/convert-db.8
man/man8/ircservices.8
sbin/ircservices
sbin/ircservices-chk
share/examples/ircservices/modules.conf
share/examples/ircservices/ircservices.conf
share/doc/html/ircservices/0.html
share/doc/html/ircservices/1.html
share/doc/html/ircservices/2.html
share/doc/html/ircservices/3.html
share/doc/html/ircservices/4.html
share/doc/html/ircservices/5.html
share/doc/html/ircservices/6.html
share/doc/html/ircservices/7.html
share/doc/html/ircservices/COPYING.html
share/doc/html/ircservices/a.html
share/doc/html/ircservices/b.html
share/doc/html/ircservices/c.html
share/doc/html/ircservices/copyright.html
share/doc/html/ircservices/faq.html
share/doc/html/ircservices/index.html
share/doc/html/ircservices/upgrade.html
share/examples/rc.d/ircservices
@dirrm share/doc/html/ircservices
@dirrm share/examples/ircservices
@dirrm lib/services/helpfiles/ircii/alias
@dirrm lib/services/helpfiles/ircii/bind
@dirrm lib/services/helpfiles/ircii/ctcp
@dirrm lib/services/helpfiles/ircii/dcc
@dirrm lib/services/helpfiles/ircii/ircii
@dirrm lib/services/helpfiles/ircii/load
@dirrm lib/services/helpfiles/ircii/note
@dirrm lib/services/helpfiles/ircii/on
@dirrm lib/services/helpfiles/ircii/set
@dirrm lib/services/helpfiles/ircii/window
@dirrm lib/services/helpfiles/ircii
@dirrm lib/services/helpfiles
@dirrm lib/services/languages
@dirrm lib/services
