# chat/ircservices/PLIST
@comment $NetBSD: PLIST,v 1.9 2009/06/14 17:37:46 joerg Exp $
lib/services/convert-db
lib/services/helpfiles/help
lib/services/languages/de
lib/services/languages/en_us
lib/services/languages/es
lib/services/languages/fr
lib/services/languages/hu
lib/services/languages/it
lib/services/languages/ja_euc
lib/services/languages/ja_sjis
lib/services/languages/nl
lib/services/languages/pt
lib/services/languages/ru
lib/services/languages/tr
man/man8/convert-db.8
man/man8/ircservices.8
sbin/ircservices
sbin/ircservices-chk
share/examples/ircservices/modules.conf
share/examples/ircservices/ircservices.conf
share/doc/html/ircservices/0.html
share/doc/html/ircservices/1.html
share/doc/html/ircservices/2.html
share/doc/html/ircservices/3.html
share/doc/html/ircservices/4.html
share/doc/html/ircservices/5.html
share/doc/html/ircservices/6.html
share/doc/html/ircservices/7.html
share/doc/html/ircservices/COPYING.html
share/doc/html/ircservices/a.html
share/doc/html/ircservices/b.html
share/doc/html/ircservices/c.html
share/doc/html/ircservices/copyright.html
share/doc/html/ircservices/faq.html
share/doc/html/ircservices/index.html
share/doc/html/ircservices/upgrade.html
share/examples/rc.d/ircservices
@unexec ${RMDIR} %D/lib/services 2>/dev/null || ${TRUE}
