# net/vsftpd/PLIST
@comment $NetBSD: PLIST,v 1.6 2006/01/13 18:12:46 wiz Exp $
libexec/vsftpd
man/man5/vsftpd.conf.5
man/man8/vsftpd.8
share/doc/vsftpd/FAQ
share/doc/vsftpd/INSTALL
share/doc/vsftpd/README
share/doc/vsftpd/README.ssl
share/doc/vsftpd/SECURITY/DESIGN
share/doc/vsftpd/SECURITY/IMPLEMENTATION
share/doc/vsftpd/SECURITY/OVERVIEW
share/doc/vsftpd/SECURITY/TRUST
share/doc/vsftpd/TUNING
share/examples/rc.d/vsftpd
share/examples/vsftpd/vsftpd.conf.default
@dirrm share/examples/vsftpd
@dirrm share/doc/vsftpd/SECURITY
@dirrm share/doc/vsftpd
