# mail/jchkmail/patches/patch-aa
$NetBSD: patch-aa,v 1.2 2007/09/04 20:40:43 joerg Exp $

--- bin/jchkmail.init.in.orig	2007-09-04 19:50:09.000000000 +0000
+++ bin/jchkmail.init.in
@@ -1,19 +1,8 @@
-#!/bin/sh
+#!@RCD_SCRIPTS_SHELL@
 #
-# chkconfig: 2345 80 30
-# description: j-chkmail - Joe's j-chkmail
-#		filtre de messagerie
-
-### BEGIN INIT INFO
-# Provides:                     j-chkmail
-# Required-Start:               $remote_fs
-# Required-Stop:
-# X-UnitedLinux-Should-Start:
-# X-UnitedLinux-Should-Stop:
-# Default-Start:                3 5
-# Default-Stop:
-# Description:                  Start j-chmail (filtre mail)
-### END INIT INFO
+# PROVIDES: jchkmail
+# BEFORE: mail
+# REQUIRE: LOGIN
 
 JCHKMAIL_HOME=@prefix@/sbin
 
