# security/courier-authlib/PLIST
@comment $NetBSD: PLIST,v 1.4 2005/05/02 20:34:04 reed Exp $
bin/courierauthconfig
bin/makedat
include/courier_auth_config.h
include/courierauth.h
include/courierauthdebug.h
include/courierauthsasl.h
include/courierauthsaslclient.h
lib/courier-authlib/libcourierauth.la
lib/courier-authlib/libcourierauthcommon.la
lib/courier-authlib/libcourierauthsasl.la
lib/courier-authlib/libcourierauthsaslclient.la
libexec/courier-authlib/authdaemond
libexec/courier-authlib/authsystem.passwd
man/man1/courierlogger.1
man/man3/auth_enumerate.3
man/man3/auth_generic.3
man/man3/auth_getoption.3
man/man3/auth_getuserinfo.3
man/man3/auth_login.3
man/man3/auth_passwd.3
man/man3/auth_sasl.3
man/man3/authlib.3
man/man8/makeuserdb.8
man/man8/userdb.8
man/man8/userdbpw.8
sbin/authdaemond
sbin/authenumerate
sbin/authtest
sbin/courierlogger
sbin/sysconftool
share/doc/courier-authlib/INSTALL.html
share/doc/courier-authlib/README.authdebug.html
share/doc/courier-authlib/README.html
share/doc/courier-authlib/README_authlib.html
share/examples/rc.d/authdaemond
