# sysutils/xentools3/PLIST
@comment $NetBSD: PLIST,v 1.1.1.1 2007/06/14 19:39:45 bouyer Exp $
bin/readnotes
bin/xc_restore
bin/xc_save
bin/xencons
bin/xenconsole
bin/xenperf
bin/xenstore-chmod
bin/xenstore-control
bin/xenstore-exists
bin/xenstore-list
bin/xenstore-ls
bin/xenstore-read
bin/xenstore-rm
bin/xenstore-write
bin/xen-detect
include/xenctrl.h
include/xenguest.h
include/xs.h
include/xs_lib.h
lib/libxenctrl.a
lib/libxenctrl.so
lib/libxenctrl.so.3.0
lib/libxenctrl.so.3.0.0
lib/libxenguest.a
lib/libxenguest.so
lib/libxenguest.so.3.0
lib/libxenguest.so.3.0.0
lib/libxenstore.a
lib/libxenstore.so
lib/libxenstore.so.3.0
lib/libxenstore.so.3.0.0
man/man1/xentop.1
man/man1/xm.1
man/man5/xend-config.sxp.5
man/man5/xmdomain.cfg.5
${PYSITELIB}/xen/__init__.py
${PYSITELIB}/xen/__init__.pyc
${PYSITELIB}/xen/lowlevel/__init__.py
${PYSITELIB}/xen/lowlevel/__init__.pyc
${PYSITELIB}/xen/lowlevel/acm.so
${PYSITELIB}/xen/lowlevel/ptsname.so
${PYSITELIB}/xen/lowlevel/xc.so
${PYSITELIB}/xen/lowlevel/xs.so
${PYSITELIB}/xen/sv/CreateDomain.py
${PYSITELIB}/xen/sv/CreateDomain.pyc
${PYSITELIB}/xen/sv/DomInfo.py
${PYSITELIB}/xen/sv/DomInfo.pyc
${PYSITELIB}/xen/sv/GenTabbed.py
${PYSITELIB}/xen/sv/GenTabbed.pyc
${PYSITELIB}/xen/sv/HTMLBase.py
${PYSITELIB}/xen/sv/HTMLBase.pyc
${PYSITELIB}/xen/sv/Main.py
${PYSITELIB}/xen/sv/Main.pyc
${PYSITELIB}/xen/sv/NodeInfo.py
${PYSITELIB}/xen/sv/NodeInfo.pyc
${PYSITELIB}/xen/sv/RestoreDomain.py
${PYSITELIB}/xen/sv/RestoreDomain.pyc
${PYSITELIB}/xen/sv/Wizard.py
${PYSITELIB}/xen/sv/Wizard.pyc
${PYSITELIB}/xen/sv/__init__.py
${PYSITELIB}/xen/sv/__init__.pyc
${PYSITELIB}/xen/sv/util.py
${PYSITELIB}/xen/sv/util.pyc
${PYSITELIB}/xen/util/Brctl.py
${PYSITELIB}/xen/util/Brctl.pyc
${PYSITELIB}/xen/util/SSHTransport.py
${PYSITELIB}/xen/util/SSHTransport.pyc
${PYSITELIB}/xen/util/__init__.py
${PYSITELIB}/xen/util/__init__.pyc
${PYSITELIB}/xen/util/asserts.py
${PYSITELIB}/xen/util/asserts.pyc
${PYSITELIB}/xen/util/auxbin.py
${PYSITELIB}/xen/util/auxbin.pyc
${PYSITELIB}/xen/util/blkif.py
${PYSITELIB}/xen/util/blkif.pyc
${PYSITELIB}/xen/util/bugtool.py
${PYSITELIB}/xen/util/bugtool.pyc
${PYSITELIB}/xen/util/diagnose.py
${PYSITELIB}/xen/util/diagnose.pyc
${PYSITELIB}/xen/util/dictio.py
${PYSITELIB}/xen/util/dictio.pyc
${PYSITELIB}/xen/util/ip.py
${PYSITELIB}/xen/util/ip.pyc
${PYSITELIB}/xen/util/mac.py
${PYSITELIB}/xen/util/mac.pyc
${PYSITELIB}/xen/util/mkdir.py
${PYSITELIB}/xen/util/mkdir.pyc
${PYSITELIB}/xen/util/pci.py
${PYSITELIB}/xen/util/pci.pyc
${PYSITELIB}/xen/util/security.py
${PYSITELIB}/xen/util/security.pyc
${PYSITELIB}/xen/util/xmlrpcclient.py
${PYSITELIB}/xen/util/xmlrpcclient.pyc
${PYSITELIB}/xen/util/xmlrpclib2.py
${PYSITELIB}/xen/util/xmlrpclib2.pyc
${PYSITELIB}/xen/util/xpopen.py
${PYSITELIB}/xen/util/xpopen.pyc
${PYSITELIB}/xen/web/SrvBase.py
${PYSITELIB}/xen/web/SrvBase.pyc
${PYSITELIB}/xen/web/SrvDir.py
${PYSITELIB}/xen/web/SrvDir.pyc
${PYSITELIB}/xen/web/__init__.py
${PYSITELIB}/xen/web/__init__.pyc
${PYSITELIB}/xen/web/connection.py
${PYSITELIB}/xen/web/connection.pyc
${PYSITELIB}/xen/web/http.py
${PYSITELIB}/xen/web/http.pyc
${PYSITELIB}/xen/web/httpserver.py
${PYSITELIB}/xen/web/httpserver.pyc
${PYSITELIB}/xen/web/protocol.py
${PYSITELIB}/xen/web/protocol.pyc
${PYSITELIB}/xen/web/resource.py
${PYSITELIB}/xen/web/resource.pyc
${PYSITELIB}/xen/web/static.py
${PYSITELIB}/xen/web/static.pyc
${PYSITELIB}/xen/web/tcp.py
${PYSITELIB}/xen/web/tcp.pyc
${PYSITELIB}/xen/web/unix.py
${PYSITELIB}/xen/web/unix.pyc
${PYSITELIB}/xen/xend/Args.py
${PYSITELIB}/xen/xend/Args.pyc
${PYSITELIB}/xen/xend/PrettyPrint.py
${PYSITELIB}/xen/xend/PrettyPrint.pyc
${PYSITELIB}/xen/xend/Vifctl.py
${PYSITELIB}/xen/xend/Vifctl.pyc
${PYSITELIB}/xen/xend/XendAPI.py
${PYSITELIB}/xen/xend/XendAPI.pyc
${PYSITELIB}/xen/xend/XendAPIConstants.py
${PYSITELIB}/xen/xend/XendAPIConstants.pyc
${PYSITELIB}/xen/xend/XendAPIStore.py
${PYSITELIB}/xen/xend/XendAPIStore.pyc
${PYSITELIB}/xen/xend/XendAPIVersion.py
${PYSITELIB}/xen/xend/XendAPIVersion.pyc
${PYSITELIB}/xen/xend/XendAuthSessions.py
${PYSITELIB}/xen/xend/XendAuthSessions.pyc
${PYSITELIB}/xen/xend/XendBase.py
${PYSITELIB}/xen/xend/XendBase.pyc
${PYSITELIB}/xen/xend/XendBootloader.py
${PYSITELIB}/xen/xend/XendBootloader.pyc
${PYSITELIB}/xen/xend/XendCheckpoint.py
${PYSITELIB}/xen/xend/XendCheckpoint.pyc
${PYSITELIB}/xen/xend/XendClient.py
${PYSITELIB}/xen/xend/XendClient.pyc
${PYSITELIB}/xen/xend/XendConfig.py
${PYSITELIB}/xen/xend/XendConfig.pyc
${PYSITELIB}/xen/xend/XendConstants.py
${PYSITELIB}/xen/xend/XendConstants.pyc
${PYSITELIB}/xen/xend/XendDevices.py
${PYSITELIB}/xen/xend/XendDevices.pyc
${PYSITELIB}/xen/xend/XendDmesg.py
${PYSITELIB}/xen/xend/XendDmesg.pyc
${PYSITELIB}/xen/xend/XendDomain.py
${PYSITELIB}/xen/xend/XendDomain.pyc
${PYSITELIB}/xen/xend/XendDomainInfo.py
${PYSITELIB}/xen/xend/XendDomainInfo.pyc
${PYSITELIB}/xen/xend/XendError.py
${PYSITELIB}/xen/xend/XendError.pyc
${PYSITELIB}/xen/xend/XendLocalStorageRepo.py
${PYSITELIB}/xen/xend/XendLocalStorageRepo.pyc
${PYSITELIB}/xen/xend/XendLogging.py
${PYSITELIB}/xen/xend/XendLogging.pyc
${PYSITELIB}/xen/xend/XendMonitor.py
${PYSITELIB}/xen/xend/XendMonitor.pyc
${PYSITELIB}/xen/xend/XendNetwork.py
${PYSITELIB}/xen/xend/XendNetwork.pyc
${PYSITELIB}/xen/xend/XendNode.py
${PYSITELIB}/xen/xend/XendNode.pyc
${PYSITELIB}/xen/xend/XendOptions.py
${PYSITELIB}/xen/xend/XendOptions.pyc
${PYSITELIB}/xen/xend/XendPBD.py
${PYSITELIB}/xen/xend/XendPBD.pyc
${PYSITELIB}/xen/xend/XendPIF.py
${PYSITELIB}/xen/xend/XendPIF.pyc
${PYSITELIB}/xen/xend/XendPIFMetrics.py
${PYSITELIB}/xen/xend/XendPIFMetrics.pyc
${PYSITELIB}/xen/xend/XendProtocol.py
${PYSITELIB}/xen/xend/XendProtocol.pyc
${PYSITELIB}/xen/xend/XendQCoWStorageRepo.py
${PYSITELIB}/xen/xend/XendQCoWStorageRepo.pyc
${PYSITELIB}/xen/xend/XendStateStore.py
${PYSITELIB}/xen/xend/XendStateStore.pyc
${PYSITELIB}/xen/xend/XendStorageRepository.py
${PYSITELIB}/xen/xend/XendStorageRepository.pyc
${PYSITELIB}/xen/xend/XendTask.py
${PYSITELIB}/xen/xend/XendTask.pyc
${PYSITELIB}/xen/xend/XendTaskManager.py
${PYSITELIB}/xen/xend/XendTaskManager.pyc
${PYSITELIB}/xen/xend/XendVDI.py
${PYSITELIB}/xen/xend/XendVDI.pyc
${PYSITELIB}/xen/xend/XendVMMetrics.py
${PYSITELIB}/xen/xend/XendVMMetrics.pyc
${PYSITELIB}/xen/xend/XendVnet.py
${PYSITELIB}/xen/xend/XendVnet.pyc
${PYSITELIB}/xen/xend/__init__.py
${PYSITELIB}/xen/xend/__init__.pyc
${PYSITELIB}/xen/xend/arch.py
${PYSITELIB}/xen/xend/arch.pyc
${PYSITELIB}/xen/xend/balloon.py
${PYSITELIB}/xen/xend/balloon.pyc
${PYSITELIB}/xen/xend/encode.py
${PYSITELIB}/xen/xend/encode.pyc
${PYSITELIB}/xen/xend/image.py
${PYSITELIB}/xen/xend/image.pyc
${PYSITELIB}/xen/xend/osdep.py
${PYSITELIB}/xen/xend/osdep.pyc
${PYSITELIB}/xen/xend/server/BlktapController.py
${PYSITELIB}/xen/xend/server/BlktapController.pyc
${PYSITELIB}/xen/xend/server/ConsoleController.py
${PYSITELIB}/xen/xend/server/ConsoleController.pyc
${PYSITELIB}/xen/xend/server/DevController.py
${PYSITELIB}/xen/xend/server/DevController.pyc
${PYSITELIB}/xen/xend/server/SSLXMLRPCServer.py
${PYSITELIB}/xen/xend/server/SSLXMLRPCServer.pyc
${PYSITELIB}/xen/xend/server/SrvDaemon.py
${PYSITELIB}/xen/xend/server/SrvDaemon.pyc
${PYSITELIB}/xen/xend/server/SrvDmesg.py
${PYSITELIB}/xen/xend/server/SrvDmesg.pyc
${PYSITELIB}/xen/xend/server/SrvDomain.py
${PYSITELIB}/xen/xend/server/SrvDomain.pyc
${PYSITELIB}/xen/xend/server/SrvDomainDir.py
${PYSITELIB}/xen/xend/server/SrvDomainDir.pyc
${PYSITELIB}/xen/xend/server/SrvNode.py
${PYSITELIB}/xen/xend/server/SrvNode.pyc
${PYSITELIB}/xen/xend/server/SrvRoot.py
${PYSITELIB}/xen/xend/server/SrvRoot.pyc
${PYSITELIB}/xen/xend/server/SrvServer.py
${PYSITELIB}/xen/xend/server/SrvServer.pyc
${PYSITELIB}/xen/xend/server/SrvVnetDir.py
${PYSITELIB}/xen/xend/server/SrvVnetDir.pyc
${PYSITELIB}/xen/xend/server/SrvXendLog.py
${PYSITELIB}/xen/xend/server/SrvXendLog.pyc
${PYSITELIB}/xen/xend/server/XMLRPCServer.py
${PYSITELIB}/xen/xend/server/XMLRPCServer.pyc
${PYSITELIB}/xen/xend/server/__init__.py
${PYSITELIB}/xen/xend/server/__init__.pyc
${PYSITELIB}/xen/xend/server/blkif.py
${PYSITELIB}/xen/xend/server/blkif.pyc
${PYSITELIB}/xen/xend/server/iopif.py
${PYSITELIB}/xen/xend/server/iopif.pyc
${PYSITELIB}/xen/xend/server/irqif.py
${PYSITELIB}/xen/xend/server/irqif.pyc
${PYSITELIB}/xen/xend/server/netif.py
${PYSITELIB}/xen/xend/server/netif.pyc
${PYSITELIB}/xen/xend/server/params.py
${PYSITELIB}/xen/xend/server/params.pyc
${PYSITELIB}/xen/xend/server/pciif.py
${PYSITELIB}/xen/xend/server/pciif.pyc
${PYSITELIB}/xen/xend/server/pciquirk.py
${PYSITELIB}/xen/xend/server/pciquirk.pyc
${PYSITELIB}/xen/xend/server/relocate.py
${PYSITELIB}/xen/xend/server/relocate.pyc
${PYSITELIB}/xen/xend/server/tests/__init__.py
${PYSITELIB}/xen/xend/server/tests/__init__.pyc
${PYSITELIB}/xen/xend/server/tests/test_controllers.py
${PYSITELIB}/xen/xend/server/tests/test_controllers.pyc
${PYSITELIB}/xen/xend/server/tpmif.py
${PYSITELIB}/xen/xend/server/tpmif.pyc
${PYSITELIB}/xen/xend/server/usbif.py
${PYSITELIB}/xen/xend/server/usbif.pyc
${PYSITELIB}/xen/xend/server/vfbif.py
${PYSITELIB}/xen/xend/server/vfbif.pyc
${PYSITELIB}/xen/xend/sxp.py
${PYSITELIB}/xen/xend/sxp.pyc
${PYSITELIB}/xen/xend/tests/__init__.py
${PYSITELIB}/xen/xend/tests/__init__.pyc
${PYSITELIB}/xen/xend/tests/test_XendConfig.py
${PYSITELIB}/xen/xend/tests/test_XendConfig.pyc
${PYSITELIB}/xen/xend/tests/test_sxp.py
${PYSITELIB}/xen/xend/tests/test_sxp.pyc
${PYSITELIB}/xen/xend/tests/test_uuid.py
${PYSITELIB}/xen/xend/tests/test_uuid.pyc
${PYSITELIB}/xen/xend/uuid.py
${PYSITELIB}/xen/xend/uuid.pyc
${PYSITELIB}/xen/xend/xenstore/__init__.py
${PYSITELIB}/xen/xend/xenstore/__init__.pyc
${PYSITELIB}/xen/xend/xenstore/tests/__init__.py
${PYSITELIB}/xen/xend/xenstore/tests/__init__.pyc
${PYSITELIB}/xen/xend/xenstore/tests/stress_xs.py
${PYSITELIB}/xen/xend/xenstore/tests/stress_xs.pyc
${PYSITELIB}/xen/xend/xenstore/xstransact.py
${PYSITELIB}/xen/xend/xenstore/xstransact.pyc
${PYSITELIB}/xen/xend/xenstore/xsutil.py
${PYSITELIB}/xen/xend/xenstore/xsutil.pyc
${PYSITELIB}/xen/xend/xenstore/xswatch.py
${PYSITELIB}/xen/xend/xenstore/xswatch.pyc
${PYSITELIB}/xen/xm/XenAPI.py
${PYSITELIB}/xen/xm/XenAPI.pyc
${PYSITELIB}/xen/xm/__init__.py
${PYSITELIB}/xen/xm/__init__.pyc
${PYSITELIB}/xen/xm/addlabel.py
${PYSITELIB}/xen/xm/addlabel.pyc
${PYSITELIB}/xen/xm/cfgbootpolicy.py
${PYSITELIB}/xen/xm/cfgbootpolicy.pyc
${PYSITELIB}/xen/xm/console.py
${PYSITELIB}/xen/xm/console.pyc
${PYSITELIB}/xen/xm/create.py
${PYSITELIB}/xen/xm/create.pyc
${PYSITELIB}/xen/xm/dry-run.py
${PYSITELIB}/xen/xm/dry-run.pyc
${PYSITELIB}/xen/xm/dumppolicy.py
${PYSITELIB}/xen/xm/dumppolicy.pyc
${PYSITELIB}/xen/xm/getlabel.py
${PYSITELIB}/xen/xm/getlabel.pyc
${PYSITELIB}/xen/xm/help.py
${PYSITELIB}/xen/xm/help.pyc
${PYSITELIB}/xen/xm/labels.py
${PYSITELIB}/xen/xm/labels.pyc
${PYSITELIB}/xen/xm/loadpolicy.py
${PYSITELIB}/xen/xm/loadpolicy.pyc
${PYSITELIB}/xen/xm/main.py
${PYSITELIB}/xen/xm/main.pyc
${PYSITELIB}/xen/xm/makepolicy.py
${PYSITELIB}/xen/xm/makepolicy.pyc
${PYSITELIB}/xen/xm/migrate.py
${PYSITELIB}/xen/xm/migrate.pyc
${PYSITELIB}/xen/xm/new.py
${PYSITELIB}/xen/xm/new.pyc
${PYSITELIB}/xen/xm/opts.py
${PYSITELIB}/xen/xm/opts.pyc
${PYSITELIB}/xen/xm/resources.py
${PYSITELIB}/xen/xm/resources.pyc
${PYSITELIB}/xen/xm/rmlabel.py
${PYSITELIB}/xen/xm/rmlabel.pyc
${PYSITELIB}/xen/xm/shutdown.py
${PYSITELIB}/xen/xm/shutdown.pyc
${PYSITELIB}/xen/xm/tests/__init__.py
${PYSITELIB}/xen/xm/tests/__init__.pyc
${PYSITELIB}/xen/xm/tests/test_create.py
${PYSITELIB}/xen/xm/tests/test_create.pyc
${PYSITELIB}/xen/xm/xenapi_create.py
${PYSITELIB}/xen/xm/xenapi_create.pyc
sbin/netfix
sbin/xen-bugtool
sbin/xenbackendd
sbin/xenconsoled
sbin/xend
sbin/xenperf
sbin/xenstored
sbin/xentop
sbin/xen-python-path
sbin/xm
share/examples/rc.d/xenbackendd
share/examples/rc.d/xend
share/examples/rc.d/xendomains
share/examples/xen/block
share/examples/xen/netbsd1
share/examples/xen/vif-bridge
share/examples/xen/vif-ip
share/examples/xen/xend-config.sxp
share/examples/xen/xmexample1
share/examples/xen/xmexample2
share/examples/xen/xmexample3
share/xen/create.dtd
@dirrm share/xen
@dirrm share/examples/xen
@dirrm ${PYSITELIB}/xen/xm/tests
@dirrm ${PYSITELIB}/xen/xm
@dirrm ${PYSITELIB}/xen/xend/xenstore/tests
@dirrm ${PYSITELIB}/xen/xend/xenstore
@dirrm ${PYSITELIB}/xen/xend/tests
@dirrm ${PYSITELIB}/xen/xend/server/tests
@dirrm ${PYSITELIB}/xen/xend/server
@dirrm ${PYSITELIB}/xen/xend
@dirrm ${PYSITELIB}/xen/web
@dirrm ${PYSITELIB}/xen/util
@dirrm ${PYSITELIB}/xen/sv
@dirrm ${PYSITELIB}/xen/lowlevel
@dirrm ${PYSITELIB}/xen
