# net/bird/PLIST
@comment $NetBSD: PLIST,v 1.1.1.1 2012/01/22 02:53:01 agc Exp $
sbin/bird${BIRD6}
sbin/birdc${BIRD6}
share/examples/bird${BIRD6}/bird${BIRD6}.conf
share/examples/rc.d/bird${BIRD6}
