# net/spread/PLIST
@comment $NetBSD: PLIST,v 1.6 2005/05/02 20:34:04 reed Exp $
bin/spflooder
bin/spmonitor
bin/sptuser
bin/spuser
include/spread/sp.h
include/spread/sp_events.h
include/spread/sp_func.h
lib/libspread.a
lib/libspread.so
lib/libspread.so.1.0
lib/libtspread.a
lib/libtspread.so
lib/libtspread.so.1.0
man/man1/spflooder.1
man/man1/spmonitor.1
man/man1/spread.1
man/man1/sptuser.1
man/man1/spuser.1
man/man3/SP_connect.3
man/man3/SP_disconnect.3
man/man3/SP_equal_group_ids.3
man/man3/SP_error.3
man/man3/SP_join.3
man/man3/SP_leave.3
man/man3/SP_multicast.3
man/man3/SP_multigroup_multicast.3
man/man3/SP_multigroup_scat_multicast.3
man/man3/SP_poll.3
man/man3/SP_receive.3
man/man3/SP_scat_multicast.3
man/man3/SP_scat_receive.3
man/man3/libsp.3
sbin/spread
share/doc/spread/LICENSE
share/doc/spread/Readme.txt
share/doc/spread/Short_Buffer_Handling.txt
share/examples/rc.d/spread
share/examples/spread/spread.access_ip.sample
share/examples/spread/spread.conf.sample
@dirrm share/examples/spread
@dirrm share/doc/spread
@dirrm include/spread
