# sysutils/apcupsd/PLIST
@comment $NetBSD: PLIST,v 1.12.18.1 2011/08/07 09:39:48 tron Exp $
${PLIST.cgi}libexec/cgi-bin/upsstats.cgi
${PLIST.cgi}libexec/cgi-bin/upsfstats.cgi
${PLIST.cgi}libexec/cgi-bin/upsimage.cgi
${PLIST.cgi}libexec/cgi-bin/multimon.cgi
sbin/apcaccess
sbin/apctest
sbin/apcupsd
man/man5/apcupsd.conf.5
man/man8/apcaccess.8
man/man8/apccontrol.8
man/man8/apctest.8
man/man8/apcupsd.8
share/doc/html/apcupsd/apcupsd.gif
share/doc/html/apcupsd/bclogo.gif
share/doc/html/apcupsd/bugs.html
share/doc/html/apcupsd/bugs.wml
share/doc/html/apcupsd/index.html
share/doc/html/apcupsd/index.wml
share/doc/html/apcupsd/license.html
share/doc/html/apcupsd/license.wml
share/doc/html/apcupsd/lists.html
share/doc/html/apcupsd/lists.wml
share/doc/html/apcupsd/mail16d.png
share/doc/html/apcupsd/menu.inc
share/doc/html/apcupsd/publishsite
share/doc/html/apcupsd/support.html
share/doc/html/apcupsd/support.wml
share/doc/html/apcupsd/systems.html
share/doc/html/apcupsd/systems.wml
share/doc/html/apcupsd/template.inc
share/examples/apcupsd/apccontrol
share/examples/apcupsd/apcupsd.conf
share/examples/apcupsd/apcupsd.css
share/examples/apcupsd/changeme
share/examples/apcupsd/commfailure
share/examples/apcupsd/commok
share/examples/apcupsd/hosts.conf
share/examples/apcupsd/multimon.conf
share/examples/apcupsd/offbattery
share/examples/apcupsd/onbattery
share/examples/rc.d/apcupsd
