# sysutils/apcupsd/PLIST
@comment $NetBSD: PLIST,v 1.4 2005/05/02 20:34:05 reed Exp $
sbin/apcaccess
sbin/apcnisd
sbin/apctest
sbin/apcupsd
man/man8/apcupsd.8
share/doc/apcupsd/README.BackUPS
share/doc/apcupsd/README.BackUPS-Pro
share/doc/apcupsd/README.NewerBackUPS-Pro
share/doc/apcupsd/README.Share-UPS
share/doc/apcupsd/README.SmartUPS-VS
share/doc/apcupsd/README.apcaccess
share/doc/apcupsd/README.autoconfig
share/doc/apcupsd/README.cable
share/doc/apcupsd/Statement.APCC
share/doc/apcupsd/apcupsd-bug-1
share/doc/apcupsd/apcupsd-bug-2
share/doc/apcupsd/minicom.txt
share/doc/apcupsd/port.gif
share/doc/html/apcupsd/apcaccess.html
share/doc/html/apcupsd/apcnisd.html
share/doc/html/apcupsd/apctest.html
share/doc/html/apcupsd/apcupsd-styles.css
share/doc/html/apcupsd/apcupsd.book
share/doc/html/apcupsd/apcupsd.gif
share/doc/html/apcupsd/apcupsd.html
share/doc/html/apcupsd/batteries.html
share/doc/html/apcupsd/bugs.html
share/doc/html/apcupsd/cables.html
share/doc/html/apcupsd/cgiprogs.html
share/doc/html/apcupsd/charging.gif
share/doc/html/apcupsd/config-examples.html
share/doc/html/apcupsd/configure.html
share/doc/html/apcupsd/data.html
share/doc/html/apcupsd/eprom.html
share/doc/html/apcupsd/events.html
share/doc/html/apcupsd/faq.html
share/doc/html/apcupsd/index.html
share/doc/html/apcupsd/install.html
share/doc/html/apcupsd/invoking.html
share/doc/html/apcupsd/kernel_config.html
share/doc/html/apcupsd/license.html
share/doc/html/apcupsd/logging.html
share/doc/html/apcupsd/master-slave.html
share/doc/html/apcupsd/multiUPS.html
share/doc/html/apcupsd/multimon.gif
share/doc/html/apcupsd/new_features_3_7_0.html
share/doc/html/apcupsd/new_features_3_8_0.html
share/doc/html/apcupsd/new_features_3_8_1.html
share/doc/html/apcupsd/new_features_3_8_2.html
share/doc/html/apcupsd/new_features_3_8_3.html
share/doc/html/apcupsd/new_features_3_8_4.html
share/doc/html/apcupsd/new_features_3_8_5.html
share/doc/html/apcupsd/oldversions.html
share/doc/html/apcupsd/onbatt.gif
share/doc/html/apcupsd/online.gif
share/doc/html/apcupsd/security.html
share/doc/html/apcupsd/shutdown.html
share/doc/html/apcupsd/status.gif
share/doc/html/apcupsd/status.html
share/doc/html/apcupsd/stopping.html
share/doc/html/apcupsd/testing.html
share/doc/html/apcupsd/thanks.gif
share/doc/html/apcupsd/thanks.html
share/doc/html/apcupsd/troubles.html
share/doc/html/apcupsd/upgrading_to_3_7_0.html
share/doc/html/apcupsd/upsbible.html
share/doc/html/apcupsd/usb.html
share/doc/html/apcupsd/win32.html
share/doc/html/apcupsd/wininstall1.gif
share/doc/html/apcupsd/wininstall2.gif
share/doc/html/apcupsd/wininstall3.gif
share/doc/html/apcupsd/wininstall4.gif
share/doc/html/apcupsd/wininstall5.gif
share/doc/html/apcupsd/wininstall6.gif
share/doc/html/apcupsd/wininstall7.gif
share/examples/apcupsd/apcupsd.master.conf
share/examples/apcupsd/apcupsd.slave.conf
share/examples/apcupsd/apccontrol
share/examples/apcupsd/changeme
share/examples/apcupsd/commfailure
share/examples/apcupsd/commok
share/examples/apcupsd/onbattery
share/examples/apcupsd/mainsback
share/examples/rc.d/apcupsd
share/locale/it/LC_MESSAGES/apcupsd.mo
@dirrm share/doc/apcupsd
@dirrm share/doc/html/apcupsd
@dirrm share/examples/apcupsd
