# sysutils/apcupsd/PLIST
@comment $NetBSD: PLIST,v 1.5 2006/11/15 21:46:25 bouyer Exp $
sbin/apcaccess
sbin/apctest
sbin/apcupsd
man/man8/apcupsd.8
man/man8/apctest.8
share/doc/html/apcupsd/apcupsd.gif
share/doc/html/apcupsd/bclogo.gif
share/doc/html/apcupsd/bugs.html
share/doc/html/apcupsd/bugs.wml
share/doc/html/apcupsd/index.html
share/doc/html/apcupsd/index.wml
share/doc/html/apcupsd/license.html
share/doc/html/apcupsd/license.wml
share/doc/html/apcupsd/lists.html
share/doc/html/apcupsd/lists.wml
share/doc/html/apcupsd/mail16d.png
share/doc/html/apcupsd/menu.inc
share/doc/html/apcupsd/publishsite
share/doc/html/apcupsd/support.html
share/doc/html/apcupsd/support.wml
share/doc/html/apcupsd/systems.html
share/doc/html/apcupsd/systems.wml
share/doc/html/apcupsd/template.inc
share/examples/apcupsd/apcupsd.master.conf
share/examples/apcupsd/apcupsd.slave.conf
share/examples/apcupsd/apccontrol
share/examples/apcupsd/changeme
share/examples/apcupsd/commfailure
share/examples/apcupsd/commok
share/examples/apcupsd/masterconnect
share/examples/apcupsd/mastertimeout
share/examples/apcupsd/onbattery
share/examples/apcupsd/offbattery
share/examples/rc.d/apcupsd
@dirrm share/doc/html/apcupsd
@dirrm share/examples/apcupsd
