# sysutils/apcupsd/PLIST
@comment $NetBSD: PLIST,v 1.1 2001/10/31 22:52:51 zuntum Exp $
sbin/apcaccess
sbin/apcupsd
sbin/powersc
etc/rc.d/apcupsd
man/man8/apcupsd.8
share/doc/apcupsd/README.BackUPS-Pro
share/doc/apcupsd/README.NEW
share/doc/apcupsd/README.NewerBackUPS-Pro
share/doc/apcupsd/README.Share-UPS
share/doc/apcupsd/README.SmartUPS-VS
share/doc/apcupsd/README.apcaccess
share/doc/apcupsd/README.autoconfig
share/doc/apcupsd/README.cable
share/doc/apcupsd/README.developers
share/doc/apcupsd/Statement.APCC
share/doc/apcupsd/minicom.apcupsd
share/doc/apcupsd/port.gif
share/examples/apcupsd.conf
share/locale/it/LC_MESSAGES/apcupsd.mo
@dirrm share/doc/apcupsd
