# net/snmptt/PLIST
@comment $NetBSD: PLIST,v 1.1.1.1 2011/04/13 10:23:47 bouyer Exp $
sbin/snmptt
sbin/snmptt-net-snmp-test
sbin/snmpttconvert
sbin/snmptthandler
sbin/snmptthandler-embedded
sbin/snmpttconvertmib
share/doc/html/snmptt/faqs.html
share/doc/html/snmptt/index.html
share/doc/html/snmptt/layout1.css
share/doc/html/snmptt/snmptt.html
share/doc/html/snmptt/snmpttconvert.html
share/doc/html/snmptt/snmpttconvertmib.html
share/examples/rc.d/snmptt
share/examples/snmptt/snmptt.ini
