# www/yaws/PLIST
@comment $NetBSD: PLIST,v 1.3 2013/11/13 15:59:22 is Exp $
bin/yaws
lib/yaws/ebin/haxe.beam
lib/yaws/ebin/json.beam
lib/yaws/ebin/json2.beam
lib/yaws/ebin/jsonrpc.beam
lib/yaws/ebin/mime_type_c.beam
lib/yaws/ebin/mime_types.beam
lib/yaws/ebin/yaws.app
lib/yaws/ebin/yaws.appup
lib/yaws/ebin/yaws.beam
lib/yaws/ebin/yaws_api.beam
lib/yaws/ebin/yaws_app.beam
lib/yaws/ebin/yaws_appmod_cgi.beam
lib/yaws/ebin/yaws_appmod_dav.beam
lib/yaws/ebin/yaws_appmod_fcgi.beam
lib/yaws/ebin/yaws_cgi.beam
lib/yaws/ebin/yaws_compile.beam
lib/yaws/ebin/yaws_config.beam
lib/yaws/ebin/yaws_ctl.beam
lib/yaws/ebin/yaws_debug.beam
lib/yaws/ebin/yaws_dime.beam
lib/yaws/ebin/yaws_exhtml.beam
lib/yaws/ebin/yaws_generated.beam
lib/yaws/ebin/yaws_html.beam
lib/yaws/ebin/yaws_jsonrpc.beam
lib/yaws/ebin/yaws_log.beam
lib/yaws/ebin/yaws_log_file_h.beam
lib/yaws/ebin/yaws_logger.beam
lib/yaws/ebin/yaws_ls.beam
lib/yaws/ebin/yaws_multipart.beam
lib/yaws/ebin/yaws_outmod.beam
lib/yaws/ebin/yaws_pam.beam
lib/yaws/ebin/yaws_revproxy.beam
lib/yaws/ebin/yaws_rpc.beam
lib/yaws/ebin/yaws_rss.beam
lib/yaws/ebin/yaws_runmod_lock.beam
lib/yaws/ebin/yaws_sendfile.beam
lib/yaws/ebin/yaws_server.beam
lib/yaws/ebin/yaws_session_server.beam
lib/yaws/ebin/yaws_shaper.beam
lib/yaws/ebin/yaws_soap12_lib.beam
lib/yaws/ebin/yaws_soap_lib.beam
lib/yaws/ebin/yaws_soap_srv.beam
lib/yaws/ebin/yaws_sse.beam
lib/yaws/ebin/yaws_stats.beam
lib/yaws/ebin/yaws_sup.beam
lib/yaws/ebin/yaws_sup_restarts.beam
lib/yaws/ebin/yaws_ticker.beam
lib/yaws/ebin/yaws_trace.beam
lib/yaws/ebin/yaws_vdir.beam
lib/yaws/ebin/yaws_websockets.beam
lib/yaws/ebin/yaws_xmlrpc.beam
lib/yaws/ebin/yaws_zlib.beam
lib/yaws/include/soap-envelope.hrl
lib/yaws/include/soap.hrl
lib/yaws/include/wsdl11soap12.hrl
lib/yaws/include/yaws.hrl
lib/yaws/include/yaws_api.hrl
lib/yaws/include/yaws_dav.hrl
lib/yaws/include/yaws_lock.hrl
lib/yaws/include/yaws_soap.hrl
lib/yaws/priv/charset.def
lib/yaws/priv/charset.def.in
lib/yaws/priv/envelope.xsd
lib/yaws/priv/epam
lib/yaws/priv/lib/.gitignore
lib/yaws/priv/lib/setuid_drv.so
lib/yaws/priv/mime.types
lib/yaws/priv/soap-envelope.xsd
lib/yaws/priv/soap.xsd
lib/yaws/priv/wsdl.xsd
lib/yaws/priv/wsdl11soap12.xsd
lib/yaws/priv/xml.xsd
man/man1/yaws.1
man/man5/yaws.conf.5
man/man5/yaws_api.5
man/man5/yaws_soap_lib.5
share/examples/rc.d/yaws
share/examples/yaws/yaws.conf
share/yaws/www/END
share/yaws/www/END2
share/yaws/www/HEAD
share/yaws/www/Makefile
share/yaws/www/TAB.inc
share/yaws/www/api.yaws
share/yaws/www/appmods.yaws
share/yaws/www/arg.yaws
share/yaws/www/arg2.yaws
share/yaws/www/bindings.yaws
share/yaws/www/cgi-bin/foo.py
share/yaws/www/cgi.yaws
share/yaws/www/code.yaws
share/yaws/www/code/Makefile
share/yaws/www/code/myappmod.beam
share/yaws/www/code/myappmod.erl
share/yaws/www/compile_layout.dia
share/yaws/www/compile_layout.png
share/yaws/www/configuration.yaws
share/yaws/www/contact.yaws
share/yaws/www/contributors.txt
share/yaws/www/cookies.yaws
share/yaws/www/doc.yaws
share/yaws/www/dynamic.yaws
share/yaws/www/embed.yaws
share/yaws/www/examples.yaws
share/yaws/www/favicon.ico
share/yaws/www/form.yaws
share/yaws/www/haxe_intro.yaws
share/yaws/www/haxe_sample.html
share/yaws/www/haxe_sample.yaws
share/yaws/www/icons/README
share/yaws/www/icons/a.gif
share/yaws/www/icons/alert.black.gif
share/yaws/www/icons/alert.red.gif
share/yaws/www/icons/apache_pb.gif
share/yaws/www/icons/back.gif
share/yaws/www/icons/ball.gray.gif
share/yaws/www/icons/ball.red.gif
share/yaws/www/icons/binary.gif
share/yaws/www/icons/binhex.gif
share/yaws/www/icons/blank.gif
share/yaws/www/icons/bomb.gif
share/yaws/www/icons/box1.gif
share/yaws/www/icons/box2.gif
share/yaws/www/icons/broken.gif
share/yaws/www/icons/burst.gif
share/yaws/www/icons/c.gif
share/yaws/www/icons/comp.blue.gif
share/yaws/www/icons/comp.gray.gif
share/yaws/www/icons/compressed.gif
share/yaws/www/icons/continued.gif
share/yaws/www/icons/dir.gif
share/yaws/www/icons/down.gif
share/yaws/www/icons/dvi.gif
share/yaws/www/icons/erl.gif
share/yaws/www/icons/f.gif
share/yaws/www/icons/folder.gif
share/yaws/www/icons/folder.open.gif
share/yaws/www/icons/folder.sec.gif
share/yaws/www/icons/forward.gif
share/yaws/www/icons/generic.gif
share/yaws/www/icons/generic.red.gif
share/yaws/www/icons/generic.sec.gif
share/yaws/www/icons/gnu-head-tiny.jpg
share/yaws/www/icons/hand.right.gif
share/yaws/www/icons/hand.up.gif
share/yaws/www/icons/hrl.gif
share/yaws/www/icons/icon.sheet.gif
share/yaws/www/icons/image1.gif
share/yaws/www/icons/image2.gif
share/yaws/www/icons/image3.gif
share/yaws/www/icons/index.gif
share/yaws/www/icons/layout.gif
share/yaws/www/icons/left.gif
share/yaws/www/icons/link.gif
share/yaws/www/icons/mailman-large.jpg
share/yaws/www/icons/mailman.jpg
share/yaws/www/icons/movie.gif
share/yaws/www/icons/p.gif
share/yaws/www/icons/patch.gif
share/yaws/www/icons/pdf.gif
share/yaws/www/icons/php4.gif
share/yaws/www/icons/pie0.gif
share/yaws/www/icons/pie1.gif
share/yaws/www/icons/pie2.gif
share/yaws/www/icons/pie3.gif
share/yaws/www/icons/pie4.gif
share/yaws/www/icons/pie5.gif
share/yaws/www/icons/pie6.gif
share/yaws/www/icons/pie7.gif
share/yaws/www/icons/pie8.gif
share/yaws/www/icons/portal.gif
share/yaws/www/icons/ps.gif
share/yaws/www/icons/quill.gif
share/yaws/www/icons/right.gif
share/yaws/www/icons/screw1.gif
share/yaws/www/icons/screw2.gif
share/yaws/www/icons/script.gif
share/yaws/www/icons/small/README.txt
share/yaws/www/icons/small/back.gif
share/yaws/www/icons/small/binary.gif
share/yaws/www/icons/small/binhex.gif
share/yaws/www/icons/small/blank.gif
share/yaws/www/icons/small/broken.gif
share/yaws/www/icons/small/burst.gif
share/yaws/www/icons/small/comp1.gif
share/yaws/www/icons/small/comp2.gif
share/yaws/www/icons/small/compressed.gif
share/yaws/www/icons/small/continued.gif
share/yaws/www/icons/small/dir.gif
share/yaws/www/icons/small/dir2.gif
share/yaws/www/icons/small/doc.gif
share/yaws/www/icons/small/forward.gif
share/yaws/www/icons/small/generic.gif
share/yaws/www/icons/small/generic2.gif
share/yaws/www/icons/small/generic3.gif
share/yaws/www/icons/small/image.gif
share/yaws/www/icons/small/image2.gif
share/yaws/www/icons/small/index.gif
share/yaws/www/icons/small/key.gif
share/yaws/www/icons/small/movie.gif
share/yaws/www/icons/small/patch.gif
share/yaws/www/icons/small/ps.gif
share/yaws/www/icons/small/rainbow.gif
share/yaws/www/icons/small/sound.gif
share/yaws/www/icons/small/sound2.gif
share/yaws/www/icons/small/tar.gif
share/yaws/www/icons/small/text.gif
share/yaws/www/icons/small/transfer.gif
share/yaws/www/icons/small/unknown.gif
share/yaws/www/icons/small/uu.gif
share/yaws/www/icons/sound1.gif
share/yaws/www/icons/sound2.gif
share/yaws/www/icons/sphere1.gif
share/yaws/www/icons/sphere2.gif
share/yaws/www/icons/tar.gif
share/yaws/www/icons/tex.gif
share/yaws/www/icons/text.gif
share/yaws/www/icons/transfer.gif
share/yaws/www/icons/unknown.gif
share/yaws/www/icons/up.gif
share/yaws/www/icons/uu.gif
share/yaws/www/icons/uuencoded.gif
share/yaws/www/icons/world1.gif
share/yaws/www/icons/world2.gif
share/yaws/www/icons/yaws.gif
share/yaws/www/icons/yawsY.gif
share/yaws/www/icons/yaws_head.gif
share/yaws/www/icons/yaws_pb.gif
share/yaws/www/icons/yaws_y.gif
share/yaws/www/index.yaws
share/yaws/www/internals.yaws
share/yaws/www/jsolait/jsolait.js
share/yaws/www/jsolait/lib/jsonrpc.js
share/yaws/www/jsolait/lib/urllib.js
share/yaws/www/json_intro.yaws
share/yaws/www/json_sample.html
share/yaws/www/json_sample.yaws
share/yaws/www/logger_mod.yaws
share/yaws/www/man.yaws
share/yaws/www/motivation.yaws
share/yaws/www/news
share/yaws/www/pcookie.yaws
share/yaws/www/post.yaws
share/yaws/www/privbind.yaws
share/yaws/www/process_tree.dia
share/yaws/www/process_tree.png
share/yaws/www/query.yaws
share/yaws/www/readcookie.yaws
share/yaws/www/readpcookie.yaws
share/yaws/www/redirect.yaws
share/yaws/www/redirect2.yaws
share/yaws/www/redirect3.yaws
share/yaws/www/redirect4.yaws
share/yaws/www/server_sent_events.html
share/yaws/www/server_sent_events.yaws
share/yaws/www/session.yaws
share/yaws/www/session1.yaws
share/yaws/www/setcookie.yaws
share/yaws/www/setpcookie.yaws
share/yaws/www/shoppingcart/Makefile
share/yaws/www/shoppingcart/buy.yaws
share/yaws/www/shoppingcart/index.yaws
share/yaws/www/shoppingcart/junk.jpg
share/yaws/www/shoppingcart/loginpost.yaws
share/yaws/www/shoppingcart/logout.yaws
share/yaws/www/shoppingcart/shopcart.beam
share/yaws/www/shoppingcart/shopcart.erl
share/yaws/www/shoppingcart/shopcart_form.yaws
share/yaws/www/shoppingcart/source.html
share/yaws/www/shoppingcart/style.css
share/yaws/www/simple.yaws
share/yaws/www/simple_ex1.yaws
share/yaws/www/simple_ex2.yaws
share/yaws/www/simple_ex3.yaws
share/yaws/www/simple_ex4.yaws
share/yaws/www/small.yaws
share/yaws/www/soap_intro.yaws
share/yaws/www/spacer.gif
share/yaws/www/ssi.yaws
share/yaws/www/ssi/dynamic.1
share/yaws/www/ssi/dynamic.2
share/yaws/www/ssi/dynamic.3
share/yaws/www/ssi_ex1
share/yaws/www/static.html
share/yaws/www/stats_ex.yaws
share/yaws/www/stil.css
share/yaws/www/stream.yaws
share/yaws/www/todo.yaws
share/yaws/www/upload.yaws
share/yaws/www/upload0.yaws
share/yaws/www/urandom.yaws
share/yaws/www/urandom_from_pid.yaws
share/yaws/www/websockets.yaws
share/yaws/www/websockets_autobahn_endpoint.yaws
share/yaws/www/websockets_example.yaws
share/yaws/www/websockets_example_endpoint.yaws
share/yaws/www/yapp_intro.yaws
share/yaws/www/yaws-1.55_to_1.56.patch
share/yaws/www/yaws.eps
share/yaws/www/ybed.erl
share/yaws/www/ybed_sup.erl
share/yaws/www/yman.yaws
@pkgdir etc/yaws
