# mail/dovecot/PLIST
@comment $NetBSD: PLIST,v 1.6 2005/05/02 20:33:59 reed Exp $
libexec/dovecot/dovecot-auth
libexec/dovecot/imap
libexec/dovecot/imap-login
libexec/dovecot/pop3
libexec/dovecot/pop3-login
sbin/dovecot
share/doc/dovecot/USE-WIKI-INSTEAD
share/doc/dovecot/auth.txt
share/doc/dovecot/configuration.txt
share/doc/dovecot/design.txt
share/doc/dovecot/index.txt
share/doc/dovecot/mail-storages.txt
share/doc/dovecot/multiaccess.txt
share/doc/dovecot/nfs.txt
share/doc/dovecot/securecoding.txt
share/examples/dovecot/dovecot-ldap.conf
share/examples/dovecot/dovecot-openssl.cnf
share/examples/dovecot/dovecot-mysql.conf
share/examples/dovecot/dovecot-pgsql.conf
share/examples/dovecot/dovecot.conf.default
share/examples/dovecot/mkcert.sh
share/examples/rc.d/dovecot
@dirrm share/examples/dovecot
@dirrm share/doc/dovecot
@dirrm libexec/dovecot
