# mail/dovecot/PLIST
@comment $NetBSD: PLIST,v 1.15 2007/03/24 14:39:43 ghen Exp $
lib/dovecot/imap/lib01_acl_plugin.so
lib/dovecot/imap/lib02_lazy_expunge_plugin.so
lib/dovecot/imap/lib10_quota_plugin.so
lib/dovecot/imap/lib11_imap_quota_plugin.la
lib/dovecot/imap/lib11_trash_plugin.so
lib/dovecot/imap/lib20_convert_plugin.so
lib/dovecot/imap/lib20_mail_log_plugin.so
lib/dovecot/imap/lib20_zlib_plugin.la
lib/dovecot/lda/lib01_acl_plugin.so
lib/dovecot/lda/lib10_quota_plugin.so
lib/dovecot/lda/lib11_trash_plugin.so
lib/dovecot/lda/lib20_convert_plugin.so
lib/dovecot/lda/lib20_mail_log_plugin.so
lib/dovecot/lib01_acl_plugin.la
lib/dovecot/lib02_lazy_expunge_plugin.la
lib/dovecot/lib10_quota_plugin.la
lib/dovecot/lib11_trash_plugin.la
lib/dovecot/lib20_convert_plugin.la
lib/dovecot/lib20_mail_log_plugin.la
lib/dovecot/pop3/lib02_lazy_expunge_plugin.so
lib/dovecot/pop3/lib10_quota_plugin.so
lib/dovecot/pop3/lib20_convert_plugin.so
lib/dovecot/pop3/lib20_mail_log_plugin.so
libexec/dovecot/checkpassword-reply
libexec/dovecot/deliver
libexec/dovecot/dict
libexec/dovecot/dovecot-auth
libexec/dovecot/gdbhelper
libexec/dovecot/idxview
libexec/dovecot/imap
libexec/dovecot/imap-login
libexec/dovecot/logview
libexec/dovecot/pop3
libexec/dovecot/pop3-login
libexec/dovecot/rawlog
libexec/dovecot/ssl-build-param
sbin/dovecot
sbin/dovecotpw
share/doc/dovecot/USE-WIKI-INSTEAD
share/doc/dovecot/auth-protocol.txt
share/doc/dovecot/securecoding.txt
share/doc/dovecot/variables.txt
share/examples/dovecot/dovecot-example.conf
share/examples/dovecot/dovecot-ldap-example.conf
share/examples/dovecot/dovecot-openssl.cnf
share/examples/dovecot/dovecot-sql-example.conf
share/examples/dovecot/mkcert.sh
share/examples/rc.d/dovecot
@dirrm share/examples/dovecot
@dirrm share/doc/dovecot
@dirrm libexec/dovecot
@dirrm lib/dovecot/pop3
@dirrm lib/dovecot/lda
@dirrm lib/dovecot/imap
@dirrm lib/dovecot
