# net/djbdns-run/PLIST
@comment $NetBSD: PLIST,v 1.4 2005/05/03 15:38:54 wiz Exp $
share/doc/djbdns-run/README.pkgsrc
share/examples/rc.d/axfrdns
share/examples/rc.d/djbdns
share/examples/rc.d/dnscache
share/examples/rc.d/rbldns
share/examples/rc.d/tinydns
@dirrm share/doc/djbdns-run
