# mail/clamav/PLIST
@comment $NetBSD: PLIST,v 1.15 2005/05/02 20:33:59 reed Exp $
bin/clamav-config
bin/clamdscan
bin/clamscan
bin/freshclam
bin/sigtool
include/clamav.h
lib/libclamav.la
lib/pkgconfig/libclamav.pc
man/man1/clamdscan.1
man/man1/clamscan.1
man/man1/freshclam.1
man/man1/sigtool.1
man/man5/clamd.conf.5
man/man5/freshclam.conf.5
${MILTER}man/man8/clamav-milter.8
man/man8/clamd.8
${MILTER}sbin/clamav-milter
sbin/clamd
share/examples/clamav/clamd.conf
share/examples/clamav/daily.cvd
share/examples/clamav/freshclam.conf
share/examples/clamav/main.cvd
share/examples/rc.d/freshclamd
share/examples/rc.d/clamd
@dirrm share/examples/clamav
