# net/speedtouch/PLIST
@comment $NetBSD: PLIST,v 1.6 2005/05/02 20:34:04 reed Exp $
share/examples/rc.d/adsl
@exec ${MKDIR} ${PKG_SYSCONFDIR}/speedtouch
bin/speedtouch-setup
bin/speedtouch-start
bin/speedtouch-stop
libdata/speedtouch/kqd6_3.012
libdata/speedtouch/zzzl_3.012
man/man1/speedtouch-setup.1
man/man1/speedtouch-start.1
man/man1/speedtouch-stop.1
man/man8/modem_run.8
man/man8/pppoa2.8
man/man8/pppoa3.8
sbin/modem_run
sbin/pppoa2
sbin/pppoa3
share/doc/speedtouch/FAQ
share/doc/speedtouch/FAQ-fr
share/doc/speedtouch/INSTALL
share/doc/speedtouch/INSTALL-de
share/doc/speedtouch/INSTALL-fr
share/doc/speedtouch/INSTALL-it
share/doc/speedtouch/INSTALL-nl
share/doc/speedtouch/ppp.conf.sample
share/examples/speedtouch/adsl.sh.sample
share/examples/speedtouch/ppp.conf.sample
share/speedtouch/boot.v123.bin
@dirrm share/speedtouch
@dirrm share/examples/speedtouch
@dirrm share/doc/speedtouch
@dirrm libdata/speedtouch
