# net/rp-pppoe/PLIST
@comment $NetBSD: PLIST,v 1.5 2009/06/14 18:09:41 joerg Exp $
man/man5/pppoe.conf.5
man/man8/pppoe-connect.8
man/man8/pppoe-server.8
man/man8/pppoe-setup.8
man/man8/pppoe-start.8
man/man8/pppoe-status.8
man/man8/pppoe-stop.8
man/man8/pppoe.8
sbin/pppoe
sbin/pppoe-connect
sbin/pppoe-server
sbin/pppoe-setup
sbin/pppoe-start
sbin/pppoe-status
sbin/pppoe-stop
share/doc/${PKGNAME}/CHANGES
share/doc/${PKGNAME}/HOW-TO-CONNECT
share/doc/${PKGNAME}/LICENSE
share/doc/${PKGNAME}/README
share/doc/${PKGNAME}/SERVPOET
share/doc/${PKGNAME}/pap-secrets
share/examples/rc.d/pppoe
share/examples/rp-pppoe/pppoe-server-options
share/examples/rp-pppoe/pppoe.conf
