# net/rp-pppoe/PLIST
@comment $NetBSD: PLIST,v 1.3 2005/05/02 20:34:03 reed Exp $
man/man5/pppoe.conf.5
man/man8/adsl-connect.8
man/man8/adsl-setup.8
man/man8/adsl-start.8
man/man8/adsl-status.8
man/man8/adsl-stop.8
man/man8/pppoe.8
man/man8/pppoe-server.8
sbin/adsl-connect
sbin/adsl-setup
sbin/adsl-start
sbin/adsl-status
sbin/adsl-stop
sbin/pppoe
sbin/pppoe-server
share/doc/rp-pppoe/CHANGES
share/doc/rp-pppoe/HOW-TO-CONNECT
share/doc/rp-pppoe/LICENSE
share/doc/rp-pppoe/README
share/doc/rp-pppoe/SERVPOET
share/doc/rp-pppoe/pap-secrets
share/examples/rc.d/pppoe
share/examples/rp-pppoe/pppoe.conf
share/examples/rp-pppoe/pppoe-server-options
@dirrm share/examples/rp-pppoe
@dirrm share/doc/rp-pppoe
