# mail/teapop/PLIST
@comment $NetBSD: PLIST,v 1.2 2002/05/07 22:49:19 hubertf Exp $
etc/rc.d/teapop
etc/teapop.passwd
libexec/teapop
man/man8/teapop.8
sbin/cronpopauth.pl
share/examples/teapop/teapop.passwd
@dirrm share/examples/teapop
