# net/arla/PLIST
@comment $NetBSD: PLIST,v 1.18 2005/05/02 20:34:01 reed Exp $
bin/ydr
bin/rxdebug
bin/rxperf
bin/amon
bin/fs
bin/tokens
bin/afstool
bin/pts
bin/udebug
bin/afsfsperf
bin/vos
bin/bos
bin/kalog
bin/arla-send-pr
bin/mk_cmds
include/atypes.h
include/fs_errors.h
include/base64.h
include/getarg.h
include/parse_bytes.h
include/parse_time.h
include/parse_units.h
include/resolve.h
include/roken-common.h
include/rtbl.h
include/xdbm.h
include/roken.h
include/lock.h
include/preempt.h
include/timer.h
include/lwp.h
include/vldb.h
include/cb.h
include/fs.h
include/pts.h
include/vldb.cs.h
include/vldb.ss.h
include/cb.cs.h
include/cb.ss.h
include/fs.cs.h
include/fs.ss.h
include/pts.cs.h
include/pts.ss.h
include/ubik.cs.h
include/ubik.ss.h
include/bos.cs.h
include/sl.h
include/ss/ss.h
include/bool.h
include/hash.h
include/heap.h
include/list.h
include/log.h
include/arlamath.h
lib/libroken.la
lib/liblwp.a
lib/librxdefclient.a
lib/librxfsserver.a
lib/librxvlserver.a
lib/librxvolserver.a
lib/libptserver.a
lib/libptclient.a
lib/libbosserver.a
lib/libbosclient.a
lib/libkaclient.a
lib/libkaserver.a
lib/librxcbclient.a
lib/librxcbserver.a
lib/librx.la
lib/libbufdir.a
lib/libarlalib.la
lib/libsl.la
lib/libss.la
lib/libarlautil.la
lib/emacs/site-lisp/send-pr.el
libexec/arlad
man/cat3/getarg.0
man/cat3/arg_printusage.0
man/man1/rxdebug.1
man/man1/tokens.1
man/man1/afstool.1
man/man1/pts.1
man/man1/udebug.1
man/man1/kalog.1
man/man1/arla-send-pr.1
man/man3/getarg.3
man/man3/arg_printusage.3
man/man3/log_log.3
man/man3/log_vlog.3
man/man3/log_open.3
man/man3/log_close.3
man/man3/log_unit_init.3
man/man3/log_unit_free.3
man/man3/log_set_mask.3
man/man3/log_get_mask.3
man/man3/log_mask2str.3
man/man3/log_set_mask_str.3
man/man4/nnpfs.4
man/man5/AliasDB.5
man/man5/CellServDB.5
man/man5/DynRootDB.5
man/man5/SuidCells.5
man/man5/ThisCell.5
man/man5/arla.conf.5
man/man8/mount_nnpfs.8
man/man8/umount_nnpfs.8
man/man8/arlad.8
man/man8/arla-cli.8
man/man8/vos.8
sbin/mount_nnpfs
sbin/umount_nnpfs
sbin/nnpfs_makedev
sbin/startarla
sbin/arla-cli
share/examples/arla/lkm.conf
share/examples/arla/CellServDB
share/examples/arla/ThisCell
share/examples/arla/fstab
share/examples/arla/SuidCells
share/examples/arla/arla.conf
share/examples/rc.d/arlad
@dirrm share/examples/arla
share/gnats/stacken
lkm/nnpfs_mod.o
