# net/arla/PLIST
@comment $NetBSD: PLIST,v 1.21 2006/08/24 17:59:04 wennmach Exp $
bin/afsfsperf
bin/afstool
${X11_SUPPORT}bin/amon
bin/arla-send-pr
bin/bos
bin/fs
bin/kalog
bin/mk_cmds
bin/pts
bin/rxdebug
bin/rxperf
bin/tokens
bin/udebug
bin/unlog
bin/vos
bin/ydr
include/arla/afs_uuid.h
include/arla/agetarg.h
include/arla/cellconfig.h
include/arla/ko.h
include/arla/part.h
include/arla/ports.h
include/arla/uae.h
include/arlamath.h
include/atypes.h
include/base64.h
include/bool.h
include/fs_errors.h
include/getarg.h
include/hash.h
include/heap.h
include/list.h
include/lock.h
include/log.h
include/lwp.h
include/parse_bytes.h
include/parse_time.h
include/parse_units.h
include/preempt.h
include/resolve.h
include/roken-common.h
include/roken.h
include/rtbl.h
include/rx/rx-new.h
include/rx/rx.h
include/rx/rx_clock.h
include/rx/rx_event.h
include/rx/rx_globs.h
include/rx/rx_kern.h
include/rx/rx_kmutex.h
include/rx/rx_locl.h
include/rx/rx_mach.h
include/rx/rx_misc.h
include/rx/rx_multi.h
include/rx/rx_null.h
include/rx/rx_pkt.h
include/rx/rx_queue.h
include/rx/rx_strm.h
include/rx/rx_trace.h
include/rx/rx_user.h
include/rx/rxgencon.h
include/rxkad/rxkad.h
include/sl.h
include/ss/ss.h
include/timer.h
include/xdbm.h
info/arla.info
info/send-pr.info
lib/emacs/site-lisp/send-pr.el
lib/libarlalib.la
lib/libarlautil.la
lib/libbosclient.la
lib/libbosserver.la
lib/libbufdir.la
lib/libkaclient.la
lib/libkaserver.la
lib/liblwp.a
lib/libptclient.la
lib/libptserver.la
lib/libroken.la
lib/librx.la
lib/librxcbclient.la
lib/librxcbserver.la
lib/librxdefclient.la
lib/librxfsserver.la
lib/librxvlserver.la
lib/librxvolserver.la
lib/libsl.la
lib/libss.la
libexec/arlad
lkm/nnpfs_mod.o
man/cat3/arg_printusage.0
man/cat3/getarg.0
man/man1/afstool.1
man/man1/arla-send-pr.1
man/man1/kalog.1
man/man1/pts.1
man/man1/rxdebug.1
man/man1/tokens.1
man/man1/udebug.1
man/man1/unlog.1
man/man3/arg_printusage.3
man/man3/getarg.3
man/man3/log_close.3
man/man3/log_get_mask.3
man/man3/log_log.3
man/man3/log_mask2str.3
man/man3/log_open.3
man/man3/log_set_mask.3
man/man3/log_set_mask_str.3
man/man3/log_unit_free.3
man/man3/log_unit_init.3
man/man3/log_vlog.3
man/man4/nnpfs.4
man/man5/AliasDB.5
man/man5/CellServDB.5
man/man5/DynRootDB.5
man/man5/SuidCells.5
man/man5/ThisCell.5
man/man5/arla.conf.5
man/man8/arla-cli.8
man/man8/arlad.8
man/man8/mount_nnpfs.8
man/man8/umount_nnpfs.8
man/man8/vos.8
sbin/arla-cli
sbin/mount_nnpfs
sbin/nnpfs_makedev
sbin/startarla
sbin/umount_nnpfs
share/examples/arla/CellServDB
share/examples/arla/DynRootDB
share/examples/arla/SuidCells
share/examples/arla/ThisCell
share/examples/arla/arla.conf
share/examples/arla/fstab
share/examples/arla/lkm.conf
@dirrm share/examples/arla
share/examples/rc.d/arlad
share/gnats/stacken
