# sysutils/gkrellm-server/PLIST
@comment $NetBSD: PLIST,v 1.1.1.1 2004/02/24 22:30:25 cube Exp $
bin/gkrellmd
etc/rc.d/gkrellmd
man/man1/gkrellmd.1
share/examples/gkrellm-server/gkrellmd.conf
@dirrm share/examples/gkrellm-server
