# sysutils/gkrellm-server/PLIST
@comment $NetBSD: PLIST,v 1.4 2005/05/02 20:34:06 reed Exp $
bin/gkrellmd
include/gkrellm2/gkrellmd.h
man/man1/gkrellmd.1
share/examples/gkrellm-server/gkrellmd.conf
share/examples/rc.d/gkrellmd
@unexec ${RMDIR} %D/include/gkrellm2 2>/dev/null || ${TRUE}
@dirrm share/examples/gkrellm-server
