# audio/musicpd/PLIST
@comment $NetBSD: PLIST,v 1.5 2009/06/14 17:32:19 joerg Exp $
bin/mpd
man/man1/mpd.1
man/man5/mpd.conf.5
share/doc/mpd/AUTHORS
share/doc/mpd/COPYING
share/doc/mpd/NEWS
share/doc/mpd/README
share/doc/mpd/TODO
share/doc/mpd/UPGRADING
share/doc/mpd/mpdconf.example
share/examples/mpd/mpd.conf
share/examples/rc.d/mpd
