# audio/musicpd/PLIST
@comment $NetBSD: PLIST,v 1.4 2009/01/20 22:52:17 drochner Exp $
bin/mpd
man/man1/mpd.1
man/man5/mpd.conf.5
share/doc/mpd/AUTHORS
share/doc/mpd/COPYING
share/doc/mpd/NEWS
share/doc/mpd/README
share/doc/mpd/TODO
share/doc/mpd/UPGRADING
share/doc/mpd/mpdconf.example
share/examples/mpd/mpd.conf
share/examples/rc.d/mpd
@dirrm share/examples/mpd
@dirrm share/doc/mpd
