# security/f-prot-antivirus6/PLIST
@comment $NetBSD: PLIST,v 1.1.1.1 2009/09/22 17:38:33 taca Exp $
bin/fpscan
${PLIST.scanmail}${PERL5_SUB_INSTALLVENDORLIB}/Anomy/HTMLCleaner.pm
${PLIST.scanmail}${PERL5_SUB_INSTALLVENDORLIB}/Anomy/Log.pm
${PLIST.scanmail}${PERL5_SUB_INSTALLVENDORLIB}/Anomy/MIMEStream.pm
${PLIST.scanmail}${PERL5_SUB_INSTALLVENDORLIB}/Anomy/Sanitizer/FProt/Client.pm
${PLIST.scanmail}${PERL5_SUB_INSTALLVENDORLIB}/Anomy/Sanitizer/FileTypes.pm
${PLIST.scanmail}${PERL5_SUB_INSTALLVENDORLIB}/Anomy/Sanitizer.pm
${PLIST.scanmail}${PERL5_SUB_INSTALLVENDORLIB}/FProt/Client.pm
${PLIST.scanmail}${FPROTDIR}/fp-milter
${PLIST.scanmail}${FPROTDIR}/fp-qmail
${PLIST.fpso}${FPROTDIR}/fp.so
${PLIST.fpso}${FPROTDIR}/fpmon
${FPROTDIR}/fpscan
${PLIST.fpscand}${FPROTDIR}/fpscand
${FPROTDIR}/fpupdate
${FPROTDIR}/product.data
${PLIST.scanmail}${FPROTDIR}/scan-mail.pl
man/man1/fpscan.1
man/man5/f-prot.conf.5
${PLIST.scanmail}man/man8/fp-milter.8
${PLIST.scanmail}man/man8/fp-qmail.8
${PLIST.fpso}man/man8/fp.so.8
${PLIST.fpso}man/man8/fpmon.8
${PLIST.fpscand}man/man8/fpscand.8
man/man8/fpupdate.8
${PLIST.scanmail}man/man8/scan-mail.pl.8
share/doc/f-prot/CHANGES
share/doc/f-prot/LICENSE
share/doc/f-prot/LICENSE-FPAV
share/doc/f-prot/LICENSES-others
share/doc/f-prot/README
share/doc/f-prot/html/add_info.html
share/doc/f-prot/html/appendix_a.html
share/doc/f-prot/html/appendix_b.html
share/doc/f-prot/html/appendix_c.html
share/doc/f-prot/html/auto_updt.html
share/doc/f-prot/html/chmstyle.css
share/doc/f-prot/html/conf_scanmail_bootup.html
share/doc/f-prot/html/fpmon.html
share/doc/f-prot/html/fpmon_conf.html
share/doc/f-prot/html/fpmon_options.html
share/doc/f-prot/html/fprots.html
share/doc/f-prot/html/fpscan.html
share/doc/f-prot/html/fpscan_options.html
share/doc/f-prot/html/fpscand.html
share/doc/f-prot/html/fpscand_conf.html
share/doc/f-prot/html/fpscand_options.html
share/doc/f-prot/html/fpupdate_options.html
share/doc/f-prot/html/images/arrow_left.gif
share/doc/f-prot/html/images/arrow_right.gif
share/doc/f-prot/html/images/banner.gif
share/doc/f-prot/html/images/dazuko.gif
share/doc/f-prot/html/images/fpupdate.gif
share/doc/f-prot/html/images/notusing_fprot_so_v6.gif
share/doc/f-prot/html/images/test_eicar.gif
share/doc/f-prot/html/images/using_fprot_so_v6.gif
share/doc/f-prot/html/index.html
share/doc/f-prot/html/inst_pro.html
share/doc/f-prot/html/loc_files.html
share/doc/f-prot/html/loglevel.html
share/doc/f-prot/html/mailscanner.html
share/doc/f-prot/html/main_feat.html
share/doc/f-prot/html/man_updt.html
share/doc/f-prot/html/milter_scan.html
share/doc/f-prot/html/per_scan.html
share/doc/f-prot/html/postfix_scan.html
share/doc/f-prot/html/prot_samba.html
share/doc/f-prot/html/qmail_scan.html
share/doc/f-prot/html/scan_with_mailscanner.html
share/doc/f-prot/html/sys_req.html
share/doc/f-prot/html/test_eicar.html
share/doc/f-prot/html/test_inst.html
share/doc/f-prot/html/updater.html
${EGDIR}/f-prot.conf
${PLIST.fpscand}share/examples/rc.d/fpscand
${PLIST.scanmail}share/examples/rc.d/scanmail
${PLIST.license}${EGDIR}/license.key
${DATADIR}/antivir.def
