# databases/mysql-server/PLIST
@comment $NetBSD: PLIST,v 1.8 2005/05/02 20:33:58 reed Exp $
bin/comp_err
bin/isamchk
bin/isamlog
bin/my_print_defaults
bin/myisam_ftdump
bin/myisamchk
bin/myisamlog
bin/myisampack
bin/mysql.server
bin/mysql_install
bin/mysql_install_db
bin/mysql_secure_installation
bin/mysql_waitpid
bin/mysqld_safe
bin/pack_isam
bin/perror
bin/replace
bin/resolve_stack_dump
bin/resolveip
lib/mysql/libdbug.a
lib/mysql/libheap.a
lib/mysql/libmerge.a
lib/mysql/libmyisam.a
lib/mysql/libmyisammrg.a
lib/mysql/libmystrings.a
lib/mysql/libmysys.a
lib/mysql/libnisam.a
lib/mysql/libvio.a
libexec/mysqld
man/man1/mysqld.1
man/man1/mysqld_multi.1
@comment man/man1/mysqld_safe.1
share/examples/rc.d/mysqld
share/mysql/charsets/Index
share/mysql/charsets/README
share/mysql/charsets/cp1251.conf
share/mysql/charsets/cp1257.conf
share/mysql/charsets/croat.conf
share/mysql/charsets/danish.conf
share/mysql/charsets/dec8.conf
share/mysql/charsets/dos.conf
share/mysql/charsets/estonia.conf
share/mysql/charsets/german1.conf
share/mysql/charsets/greek.conf
share/mysql/charsets/hebrew.conf
share/mysql/charsets/hp8.conf
share/mysql/charsets/hungarian.conf
share/mysql/charsets/koi8_ru.conf
share/mysql/charsets/koi8_ukr.conf
share/mysql/charsets/latin1.conf
share/mysql/charsets/latin2.conf
share/mysql/charsets/latin5.conf
share/mysql/charsets/swe7.conf
share/mysql/charsets/usa7.conf
share/mysql/charsets/win1250.conf
share/mysql/charsets/win1251.conf
share/mysql/charsets/win1251ukr.conf
share/mysql/czech/errmsg.sys
share/mysql/czech/errmsg.txt
share/mysql/danish/errmsg.sys
share/mysql/danish/errmsg.txt
share/mysql/dutch/errmsg.sys
share/mysql/dutch/errmsg.txt
share/mysql/english/errmsg.sys
share/mysql/english/errmsg.txt
share/mysql/estonian/errmsg.sys
share/mysql/estonian/errmsg.txt
share/mysql/french/errmsg.sys
share/mysql/french/errmsg.txt
share/mysql/german/errmsg.sys
share/mysql/german/errmsg.txt
share/mysql/greek/errmsg.sys
share/mysql/greek/errmsg.txt
share/mysql/hungarian/errmsg.sys
share/mysql/hungarian/errmsg.txt
share/mysql/italian/errmsg.sys
share/mysql/italian/errmsg.txt
share/mysql/japanese/errmsg.sys
share/mysql/japanese/errmsg.txt
share/mysql/korean/errmsg.sys
share/mysql/korean/errmsg.txt
share/mysql/mi_test_all
share/mysql/mi_test_all.res
share/mysql/mysql.server
share/mysql/norwegian-ny/errmsg.sys
share/mysql/norwegian-ny/errmsg.txt
share/mysql/norwegian/errmsg.sys
share/mysql/norwegian/errmsg.txt
share/mysql/polish/errmsg.sys
share/mysql/polish/errmsg.txt
share/mysql/portuguese/errmsg.sys
share/mysql/portuguese/errmsg.txt
share/mysql/romanian/errmsg.sys
share/mysql/romanian/errmsg.txt
share/mysql/russian/errmsg.sys
share/mysql/russian/errmsg.txt
share/mysql/slovak/errmsg.sys
share/mysql/slovak/errmsg.txt
share/mysql/spanish/errmsg.sys
share/mysql/spanish/errmsg.txt
share/mysql/swedish/errmsg.sys
share/mysql/swedish/errmsg.txt
share/mysql/ukrainian/errmsg.sys
share/mysql/ukrainian/errmsg.txt
@dirrm share/mysql/ukrainian
@dirrm share/mysql/swedish
@dirrm share/mysql/spanish
@dirrm share/mysql/slovak
@dirrm share/mysql/russian
@dirrm share/mysql/romanian
@dirrm share/mysql/portuguese
@dirrm share/mysql/polish
@dirrm share/mysql/norwegian-ny
@dirrm share/mysql/norwegian
@dirrm share/mysql/korean
@dirrm share/mysql/japanese
@dirrm share/mysql/italian
@dirrm share/mysql/hungarian
@dirrm share/mysql/greek
@dirrm share/mysql/german
@dirrm share/mysql/french
@dirrm share/mysql/estonian
@dirrm share/mysql/english
@dirrm share/mysql/dutch
@dirrm share/mysql/danish
@dirrm share/mysql/czech
@dirrm share/mysql/charsets
