# inputmethod/canna-server-bin/PLIST
@comment $NetBSD: PLIST,v 1.1.1.1.2.2 2002/06/23 18:47:26 jlam Exp $
bin/canlisp
bin/cannacheck
bin/cannakill
bin/cannastat
bin/catdic
bin/crfreq
bin/crxdic
bin/crxgram
bin/cshost
bin/ctow
bin/dicar
bin/dpbindic
bin/dpromdic
bin/dpxdic
bin/forcpp
bin/forsort
bin/kpdic
bin/mergeword
bin/mkbindic
bin/mkromdic
bin/splitword
bin/wtoc
etc/rc.d/canna
man/ja_JP.EUC/man1/addwords.1
man/ja_JP.EUC/man1/cannacheck.1
man/ja_JP.EUC/man1/cannakill.1
man/ja_JP.EUC/man1/cannaserver.1
man/ja_JP.EUC/man1/cannastat.1
man/ja_JP.EUC/man1/catdic.1
man/ja_JP.EUC/man1/cpdic.1
man/ja_JP.EUC/man1/cshost.1
man/ja_JP.EUC/man1/ctow.1
man/ja_JP.EUC/man1/delwords.1
man/ja_JP.EUC/man1/dicar.1
man/ja_JP.EUC/man1/dpbindic.1
man/ja_JP.EUC/man1/dpromdic.1
man/ja_JP.EUC/man1/lsdic.1
man/ja_JP.EUC/man1/mkbindic.1
man/ja_JP.EUC/man1/mkdic.1
man/ja_JP.EUC/man1/mkromdic.1
man/ja_JP.EUC/man1/mvdic.1
man/ja_JP.EUC/man1/rmdic.1
man/ja_JP.EUC/man1/splitword.1
man/ja_JP.EUC/man1/wtoc.1
sbin/cannaserver
@exec ${MKDIR} ${CANNA_SPOOL}
@exec ${CHOWN} ${CANNAOWNER}:${CANNAGROUP} ${CANNA_SPOOL}
@exec ${CHMOD} ${CANNA_MODE} ${CANNA_SPOOL}
@unexec ${RMDIR} ${CANNA_SPOOL}
