# net/dtcp/PLIST
@comment $NetBSD: PLIST,v 1.3 2005/05/02 20:34:02 reed Exp $
man/cat8/dtcpauth.0
man/cat8/dtcpc.0
man/cat8/dtcps.0
man/man8/dtcpauth.8
man/man8/dtcpc.8
man/man8/dtcps.8
sbin/dtcpauth
sbin/dtcpc
sbin/dtcps
share/examples/rc.d/dtcps
share/examples/rc.d/dtcpc
