# mail/sqlgrey/PLIST
@comment $NetBSD: PLIST,v 1.2 2009/06/14 18:04:43 joerg Exp $
bin/sqlgrey-logstats.pl
man/man1/sqlgrey.1
sbin/sqlgrey
sbin/update_sqlgrey_config
share/examples/rc.d/sqlgrey
share/examples/sqlgrey/README
share/examples/sqlgrey/clients_fqdn_whitelist
share/examples/sqlgrey/clients_ip_whitelist
share/examples/sqlgrey/dyn_fqdn.regexp
share/examples/sqlgrey/smtp_server.regexp
share/examples/sqlgrey/sqlgrey.conf
