# net/smokeping/PLIST
@comment $NetBSD: PLIST,v 1.12 2012/02/02 18:54:34 bouyer Exp $
bin/smokeinfo
bin/smokeping
bin/smokeping_cgi
bin/tSmoke
lib/smokeping/BER.pm
lib/smokeping/SNMP_Session.pm
lib/smokeping/SNMP_util.pm
lib/smokeping/Smokeping.pm
lib/smokeping/Smokeping/Colorspace.pm
lib/smokeping/Smokeping/Config.pm
lib/smokeping/Smokeping/Examples.pm
lib/smokeping/Smokeping/Graphs.pm
lib/smokeping/Smokeping/Info.pm
lib/smokeping/Smokeping/Master.pm
lib/smokeping/Smokeping/RRDhelpers.pm
lib/smokeping/Smokeping/RRDtools.pm
lib/smokeping/Smokeping/Slave.pm
lib/smokeping/Smokeping/ciscoRttMonMIB.pm
lib/smokeping/Smokeping/matchers/Avgratio.pm
lib/smokeping/Smokeping/matchers/CheckLatency.pm
lib/smokeping/Smokeping/matchers/CheckLoss.pm
lib/smokeping/Smokeping/matchers/ExpLoss.pm
lib/smokeping/Smokeping/matchers/Median.pm
lib/smokeping/Smokeping/matchers/Medratio.pm
lib/smokeping/Smokeping/matchers/base.pm
lib/smokeping/Smokeping/pingMIB.pm
lib/smokeping/Smokeping/probes/AnotherDNS.pm
lib/smokeping/Smokeping/probes/AnotherSSH.pm
lib/smokeping/Smokeping/probes/CiscoRTTMonDNS.pm
lib/smokeping/Smokeping/probes/CiscoRTTMonEchoICMP.pm
lib/smokeping/Smokeping/probes/CiscoRTTMonTcpConnect.pm
lib/smokeping/Smokeping/probes/Curl.pm
lib/smokeping/Smokeping/probes/DNS.pm
lib/smokeping/Smokeping/probes/DismanPing.pm
lib/smokeping/Smokeping/probes/EchoPing.pm
lib/smokeping/Smokeping/probes/EchoPingChargen.pm
lib/smokeping/Smokeping/probes/EchoPingDNS.pm
lib/smokeping/Smokeping/probes/EchoPingDiscard.pm
lib/smokeping/Smokeping/probes/EchoPingHttp.pm
lib/smokeping/Smokeping/probes/EchoPingHttps.pm
lib/smokeping/Smokeping/probes/EchoPingIcp.pm
lib/smokeping/Smokeping/probes/EchoPingLDAP.pm
lib/smokeping/Smokeping/probes/EchoPingPlugin.pm
lib/smokeping/Smokeping/probes/EchoPingSmtp.pm
lib/smokeping/Smokeping/probes/EchoPingWhois.pm
lib/smokeping/Smokeping/probes/FPing.pm
lib/smokeping/Smokeping/probes/FPing6.pm
lib/smokeping/Smokeping/probes/FTPtransfer.pm
lib/smokeping/Smokeping/probes/IOSPing.pm
lib/smokeping/Smokeping/probes/LDAP.pm
lib/smokeping/Smokeping/probes/NFSping.pm
lib/smokeping/Smokeping/probes/OpenSSHJunOSPing.pm
lib/smokeping/Smokeping/probes/Radius.pm
lib/smokeping/Smokeping/probes/RemoteFPing.pm
lib/smokeping/Smokeping/probes/SSH.pm
lib/smokeping/Smokeping/probes/SipSak.pm
lib/smokeping/Smokeping/probes/TCPPing.pm
lib/smokeping/Smokeping/probes/TacacsPlus.pm
lib/smokeping/Smokeping/probes/TelnetIOSPing.pm
lib/smokeping/Smokeping/probes/TelnetJunOSPing.pm
lib/smokeping/Smokeping/probes/WebProxyFilter.pm
lib/smokeping/Smokeping/probes/base.pm
lib/smokeping/Smokeping/probes/basefork.pm
lib/smokeping/Smokeping/probes/basevars.pm
lib/smokeping/Smokeping/probes/passwordchecker.pm
lib/smokeping/Smokeping/probes/skel.pm
lib/smokeping/Smokeping/sorters/Loss.pm
lib/smokeping/Smokeping/sorters/Max.pm
lib/smokeping/Smokeping/sorters/Median.pm
lib/smokeping/Smokeping/sorters/StdDev.pm
lib/smokeping/Smokeping/sorters/base.pm
libexec/cgi-bin/smokeping.cgi
man/man1/smokeping.1
man/man1/smokeping_cgi.1
man/man1/tSmoke.1
man/man3/Smokeping_matchers_Avgratio.3
man/man3/Smokeping_matchers_CheckLatency.3
man/man3/Smokeping_matchers_CheckLoss.3
man/man3/Smokeping_matchers_ExpLoss.3
man/man3/Smokeping_matchers_Median.3
man/man3/Smokeping_matchers_Medratio.3
man/man3/Smokeping_matchers_base.3
man/man3/Smokeping_probes_AnotherDNS.3
man/man3/Smokeping_probes_AnotherSSH.3
man/man3/Smokeping_probes_CiscoRTTMonDNS.3
man/man3/Smokeping_probes_CiscoRTTMonEchoICMP.3
man/man3/Smokeping_probes_CiscoRTTMonTcpConnect.3
man/man3/Smokeping_probes_Curl.3
man/man3/Smokeping_probes_DNS.3
man/man3/Smokeping_probes_DismanPing.3
man/man3/Smokeping_probes_EchoPing.3
man/man3/Smokeping_probes_EchoPingChargen.3
man/man3/Smokeping_probes_EchoPingDNS.3
man/man3/Smokeping_probes_EchoPingDiscard.3
man/man3/Smokeping_probes_EchoPingHttp.3
man/man3/Smokeping_probes_EchoPingHttps.3
man/man3/Smokeping_probes_EchoPingIcp.3
man/man3/Smokeping_probes_EchoPingLDAP.3
man/man3/Smokeping_probes_EchoPingPlugin.3
man/man3/Smokeping_probes_EchoPingSmtp.3
man/man3/Smokeping_probes_EchoPingWhois.3
man/man3/Smokeping_probes_FPing.3
man/man3/Smokeping_probes_FPing6.3
man/man3/Smokeping_probes_FTPtransfer.3
man/man3/Smokeping_probes_IOSPing.3
man/man3/Smokeping_probes_LDAP.3
man/man3/Smokeping_probes_NFSping.3
man/man3/Smokeping_probes_OpenSSHJunOSPing.3
man/man3/Smokeping_probes_Radius.3
man/man3/Smokeping_probes_RemoteFPing.3
man/man3/Smokeping_probes_SSH.3
man/man3/Smokeping_probes_SipSak.3
man/man3/Smokeping_probes_TCPPing.3
man/man3/Smokeping_probes_TacacsPlus.3
man/man3/Smokeping_probes_TelnetIOSPing.3
man/man3/Smokeping_probes_TelnetJunOSPing.3
man/man3/Smokeping_probes_WebProxyFilter.3
man/man3/Smokeping_probes_base.3
man/man3/Smokeping_probes_basefork.3
man/man3/Smokeping_probes_basevars.3
man/man3/Smokeping_probes_passwordchecker.3
man/man3/Smokeping_probes_skel.3
man/man3/Smokeping_sorters_Loss.3
man/man3/Smokeping_sorters_Max.3
man/man3/Smokeping_sorters_Median.3
man/man3/Smokeping_sorters_StdDev.3
man/man3/Smokeping_sorters_base.3
man/man5/smokeping_config.5
man/man5/smokeping_examples.5
man/man7/smokeping_extend.7
man/man7/smokeping_install.7
man/man7/smokeping_master_slave.7
man/man7/smokeping_upgrade.7
share/doc/smokeping/CHANGES
share/doc/smokeping/CONTRIBUTORS
share/doc/smokeping/COPYRIGHT
share/doc/smokeping/LICENSE
share/doc/smokeping/README
share/doc/smokeping/TODO
share/examples/rc.d/smokeping
share/examples/smokeping/basepage.html.dist
share/examples/smokeping/config.dist
share/examples/smokeping/config.echoping
share/examples/smokeping/config.fping-instances
share/examples/smokeping/config.multiple-probes
share/examples/smokeping/config.simple
share/examples/smokeping/config.targetvars-with-Curl
share/examples/smokeping/config.template
share/examples/smokeping/smokemail.dist
share/examples/smokeping/smokeping_secrets.dist
share/examples/smokeping/smokeping.conf
share/examples/smokeping/tmail.dist
share/httpd/htdocs/cropper/cropper.css
share/httpd/htdocs/cropper/cropper.js
share/httpd/htdocs/cropper/cropper.uncompressed.js
share/httpd/htdocs/cropper/lib/builder.js
share/httpd/htdocs/cropper/lib/controls.js
share/httpd/htdocs/cropper/lib/dragdrop.js
share/httpd/htdocs/cropper/lib/effects.js
share/httpd/htdocs/cropper/lib/prototype.js
share/httpd/htdocs/cropper/lib/scriptaculous.js
share/httpd/htdocs/cropper/lib/slider.js
share/httpd/htdocs/cropper/lib/sound.js
share/httpd/htdocs/cropper/lib/unittest.js
share/httpd/htdocs/cropper/licence.txt
share/httpd/htdocs/cropper/marqueeHoriz.gif
share/httpd/htdocs/cropper/marqueeVert.gif
share/httpd/htdocs/cropper/smokeping-zoom.js
