# net/smokeping/PLIST
@comment $NetBSD: PLIST,v 1.5 2005/05/02 20:34:04 reed Exp $
bin/smokeping
lib/smokeping/BER.pm
lib/smokeping/ISG/ParseConfig.pm
lib/smokeping/SNMP_Session.pm
lib/smokeping/SNMP_util.pm
lib/smokeping/Smokeping.pm
lib/smokeping/ciscoRttMonMIB.pm
lib/smokeping/matchers/avgratio.pm
lib/smokeping/matchers/base.pm
lib/smokeping/matchers/median.pm
lib/smokeping/probes/AnotherDNS.pm
lib/smokeping/probes/AnotherSSH.pm
lib/smokeping/probes/CiscoRTTMonDNS.pm
lib/smokeping/probes/CiscoRTTMonEchoICMP.pm
lib/smokeping/probes/CiscoRTTMonTcpConnect.pm
lib/smokeping/probes/Curl.pm
lib/smokeping/probes/DNS.pm
lib/smokeping/probes/EchoPing.pm
lib/smokeping/probes/EchoPingChargen.pm
lib/smokeping/probes/EchoPingDiscard.pm
lib/smokeping/probes/EchoPingHttp.pm
lib/smokeping/probes/EchoPingHttps.pm
lib/smokeping/probes/EchoPingIcp.pm
lib/smokeping/probes/EchoPingSmtp.pm
lib/smokeping/probes/FPing.pm
lib/smokeping/probes/FPing6.pm
lib/smokeping/probes/IOSPing.pm
lib/smokeping/probes/LDAP.pm
lib/smokeping/probes/Radius.pm
lib/smokeping/probes/RemoteFPing.pm
lib/smokeping/probes/SSH.pm
lib/smokeping/probes/base.pm
lib/smokeping/probes/basefork.pm
lib/smokeping/probes/basevars.pm
lib/smokeping/probes/passwordchecker.pm
lib/smokeping/probes/telnetIOSPing.pm
libexec/cgi-bin/smokeping.cgi
share/doc/smokeping/CHANGES
share/doc/smokeping/CONTRIBUTORS
share/doc/smokeping/COPYING
share/doc/smokeping/COPYRIGHT
share/doc/smokeping/ParseConfig.pm.txt
share/doc/smokeping/README
share/doc/smokeping/Smokeping.pm.txt
share/doc/smokeping/TODO
share/doc/smokeping/smokeping.cgi.txt
share/doc/smokeping/smokeping.txt
share/doc/smokeping/smokeping_config.txt
share/doc/smokeping/smokeping_install.txt
share/doc/smokeping/matchers/avgratio.pm.txt
share/doc/smokeping/matchers/base.pm.txt
share/doc/smokeping/matchers/median.pm.txt
share/doc/smokeping/probes/AnotherDNS.pm.txt
share/doc/smokeping/probes/AnotherSSH.pm.txt
share/doc/smokeping/probes/CiscoRTTMonDNS.pm.txt
share/doc/smokeping/probes/CiscoRTTMonEchoICMP.pm.txt
share/doc/smokeping/probes/CiscoRTTMonTcpConnect.pm.txt
share/doc/smokeping/probes/Curl.pm.txt
share/doc/smokeping/probes/DNS.pm.txt
share/doc/smokeping/probes/EchoPing.pm.txt
share/doc/smokeping/probes/EchoPingChargen.pm.txt
share/doc/smokeping/probes/EchoPingDiscard.pm.txt
share/doc/smokeping/probes/EchoPingHttp.pm.txt
share/doc/smokeping/probes/EchoPingHttps.pm.txt
share/doc/smokeping/probes/EchoPingIcp.pm.txt
share/doc/smokeping/probes/EchoPingSmtp.pm.txt
share/doc/smokeping/probes/FPing.pm.txt
share/doc/smokeping/probes/FPing6.pm.txt
share/doc/smokeping/probes/IOSPing.pm.txt
share/doc/smokeping/probes/LDAP.pm.txt
share/doc/smokeping/probes/Radius.pm.txt
share/doc/smokeping/probes/RemoteFPing.pm.txt
share/doc/smokeping/probes/SSH.pm.txt
share/doc/smokeping/probes/base.pm.txt
share/doc/smokeping/probes/basefork.pm.txt
share/doc/smokeping/probes/basevars.pm.txt
share/doc/smokeping/probes/passwordchecker.pm.txt
share/doc/smokeping/probes/telnetIOSPing.pm.txt
share/examples/rc.d/smokeping
share/examples/smokeping/basepage.html.dist
share/examples/smokeping/config-echoping.dist
share/examples/smokeping/config.dist
share/examples/smokeping/smokemail.dist
@dirrm share/examples/smokeping
@dirrm share/doc/smokeping/matchers
@dirrm share/doc/smokeping/probes
@dirrm share/doc/smokeping
@dirrm lib/smokeping/probes
@dirrm lib/smokeping/matchers
@dirrm lib/smokeping/ISG
@dirrm lib/smokeping
