# net/smokeping/PLIST
@comment $NetBSD: PLIST,v 1.7 2006/08/09 21:48:53 adrianp Exp $
bin/smokeping
bin/tSmoke
lib/smokeping/BER.pm
lib/smokeping/Config/Grammar.pm
lib/smokeping/SNMP_Session.pm
lib/smokeping/SNMP_util.pm
lib/smokeping/Smokeping.pm
lib/smokeping/Smokeping/Examples.pm
lib/smokeping/Smokeping/RRDtools.pm
lib/smokeping/Smokeping/ciscoRttMonMIB.pm
lib/smokeping/Smokeping/matchers/Avgratio.pm
lib/smokeping/Smokeping/matchers/Median.pm
lib/smokeping/Smokeping/matchers/base.pm
lib/smokeping/Smokeping/matchers/CheckLatency.pm
lib/smokeping/Smokeping/matchers/CheckLoss.pm
lib/smokeping/Smokeping/probes/AnotherDNS.pm
lib/smokeping/Smokeping/probes/AnotherSSH.pm
lib/smokeping/Smokeping/probes/CiscoRTTMonDNS.pm
lib/smokeping/Smokeping/probes/CiscoRTTMonEchoICMP.pm
lib/smokeping/Smokeping/probes/CiscoRTTMonTcpConnect.pm
lib/smokeping/Smokeping/probes/Curl.pm
lib/smokeping/Smokeping/probes/DNS.pm
lib/smokeping/Smokeping/probes/EchoPing.pm
lib/smokeping/Smokeping/probes/EchoPingChargen.pm
lib/smokeping/Smokeping/probes/EchoPingDiscard.pm
lib/smokeping/Smokeping/probes/EchoPingHttp.pm
lib/smokeping/Smokeping/probes/EchoPingHttps.pm
lib/smokeping/Smokeping/probes/EchoPingIcp.pm
lib/smokeping/Smokeping/probes/EchoPingSmtp.pm
lib/smokeping/Smokeping/probes/FPing.pm
lib/smokeping/Smokeping/probes/FPing6.pm
lib/smokeping/Smokeping/probes/IOSPing.pm
lib/smokeping/Smokeping/probes/LDAP.pm
lib/smokeping/Smokeping/probes/Radius.pm
lib/smokeping/Smokeping/probes/RemoteFPing.pm
lib/smokeping/Smokeping/probes/SSH.pm
lib/smokeping/Smokeping/probes/TacacsPlus.pm
lib/smokeping/Smokeping/probes/TelnetIOSPing.pm
lib/smokeping/Smokeping/probes/base.pm
lib/smokeping/Smokeping/probes/basefork.pm
lib/smokeping/Smokeping/probes/basevars.pm
lib/smokeping/Smokeping/probes/passwordchecker.pm
lib/smokeping/Smokeping/probes/skel.pm
libexec/cgi-bin/smokeping.cgi
share/doc/smokeping/CHANGES
share/doc/smokeping/CONTRIBUTORS
share/doc/smokeping/COPYING
share/doc/smokeping/COPYRIGHT
share/doc/smokeping/README
share/doc/smokeping/Smokeping.txt
share/doc/smokeping/TODO
share/doc/smokeping/smokeping.cgi.txt
share/doc/smokeping/smokeping.txt
share/doc/smokeping/smokeping_config.txt
share/doc/smokeping/smokeping_examples.txt
share/doc/smokeping/smokeping_extend.txt
share/doc/smokeping/smokeping_install.txt
share/doc/smokeping/smokeping_upgrade.txt
share/doc/smokeping/tSmoke.txt
share/examples/rc.d/smokeping
share/examples/smokeping/basepage.html.dist
share/examples/smokeping/config.dist
share/examples/smokeping/config.echoping
share/examples/smokeping/config.fping-instances
share/examples/smokeping/config.multiple-probes
share/examples/smokeping/config.simple
share/examples/smokeping/config.targetvars-with-Curl
share/examples/smokeping/config.template
share/examples/smokeping/smokemail.dist
share/examples/smokeping/tmail.dist
@dirrm share/examples/smokeping
@dirrm share/doc/smokeping
@dirrm lib/smokeping/Smokeping/probes
@dirrm lib/smokeping/Smokeping/matchers
@dirrm lib/smokeping/Smokeping
@dirrm lib/smokeping/Config
@dirrm lib/smokeping
