# security/cyrus-saslauthd/PLIST
@comment $NetBSD: PLIST,v 1.2 2004/01/10 01:27:00 jlam Exp $
etc/rc.d/saslauthd
man/man8/saslauthd.8
sbin/saslauthd
