# security/smtpd/PLIST
@comment $NetBSD: PLIST,v 1.3 2009/06/14 18:13:39 joerg Exp $
sbin/smtpd
sbin/smtpfwdd
share/doc/smtpd/INSTALL
share/doc/smtpd/README
share/examples/smtpd/antirelay_check_rules_example
share/examples/smtpd/antispam_check_rules.example
share/examples/smtpd/smtpd_check_rules.example
share/examples/rc.d/smtpfwdd
