# mail/spamass-milter/PLIST
@comment $NetBSD: PLIST,v 1.4 2005/05/02 20:34:01 reed Exp $
man/man1/spamass-milter.1
sbin/spamass-milter
share/examples/rc.d/spamass-milter
