# mail/postgrey/PLIST
@comment $NetBSD: PLIST,v 1.3 2006/02/19 01:01:39 wiz Exp $
sbin/postgrey
share/doc/postgrey/COPYING
share/doc/postgrey/Changes
share/doc/postgrey/README
share/examples/postgrey/postgrey_whitelist_clients
share/examples/postgrey/postgrey_whitelist_recipients
share/examples/rc.d/postgrey
@dirrm share/examples/postgrey
@dirrm share/doc/postgrey
