# www/squid/PLIST
@comment $NetBSD: PLIST,v 1.10 2002/11/25 14:46:14 taca Exp $
bin/squidclient
etc/rc.d/squid
libexec/cachemgr.cgi
man/man8/squid.8
share/examples/squid/mime.conf
share/examples/squid/squid.conf
sbin/RunAccel
sbin/RunCache
sbin/squid
@dirrm share/examples/squid
@exec ${MKDIR} ${PKG_SYSCONFDIR}
@unexec rmdir ${PKG_SYSCONFDIR} 2>/dev/null || true
