# www/squid/PLIST
@comment $NetBSD: PLIST,v 1.1 2001/11/01 02:18:53 zuntum Exp $
etc/rc.d/squid
etc/squid/mime.conf.default
etc/squid/squid.conf.default
@unexec rmdir %D/etc/squid 2>/dev/null || true
libexec/cachemgr.cgi
sbin/RunAccel
sbin/RunCache
sbin/client
sbin/squid
