# www/squid/PLIST
@comment $NetBSD: PLIST,v 1.12 2003/03/25 17:44:22 taca Exp $
bin/squidclient
etc/rc.d/squid
libexec/cachemgr.cgi
man/man8/squid.8
share/doc/squid/ChangeLog
share/doc/squid/RELEASENOTES.html
share/doc/squid/debug-sections.txt
share/examples/squid/mime.conf
share/examples/squid/squid.conf
sbin/RunAccel
sbin/RunCache
sbin/squid
@dirrm share/examples/squid
@exec ${MKDIR} ${PKG_SYSCONFDIR}
@unexec rmdir ${PKG_SYSCONFDIR} 2>/dev/null || ${TRUE}
