# www/squid/PLIST
@comment $NetBSD: PLIST,v 1.4 2001/12/12 17:06:19 taca Exp $
etc/rc.d/squid
etc/squid/mime.conf.default
etc/squid/squid.conf.default
@unexec ${RMDIR} %D/etc/squid 2>/dev/null || ${TRUE}
libexec/cachemgr.cgi
sbin/RunAccel
sbin/RunCache
sbin/client
sbin/squid
