# www/squid/PLIST
@comment $NetBSD: PLIST,v 1.7 2002/09/15 14:39:51 grant Exp $
etc/rc.d/squid
@unexec ${RMDIR} %D/etc/squid 2>/dev/null || ${TRUE}
libexec/cachemgr.cgi
share/examples/squid/mime.conf
share/examples/squid/squid.conf
sbin/RunAccel
sbin/RunCache
sbin/client
sbin/squid
@dirrm share/examples/squid
@exec ${MKDIR} ${PKG_SYSCONFDIR}
@unexec rmdir ${PKG_SYSCONFDIR} 2>/dev/null || true
