# www/squid/PLIST
@comment $NetBSD: PLIST,v 1.6 2002/06/26 10:30:08 seb Exp $
etc/rc.d/squid
@unexec ${RMDIR} %D/etc/squid 2>/dev/null || ${TRUE}
libexec/cachemgr.cgi
share/examples/squid/mime.conf
share/examples/squid/squid.conf
sbin/RunAccel
sbin/RunCache
sbin/client
sbin/squid
@dirrm share/examples/squid
@exec ${MKDIR} ${SQUID_SYSCONFDIR}
@unexec rmdir ${SQUID_SYSCONFDIR} 2>/dev/null || true
