# www/squid/PLIST
@comment $NetBSD: PLIST,v 1.17 2005/05/02 20:34:07 reed Exp $
bin/squidclient
libexec/cachemgr.cgi
man/man8/squid.8
sbin/RunAccel
sbin/RunCache
sbin/squid
share/doc/squid/ChangeLog
share/doc/squid/RELEASENOTES.html
share/doc/squid/debug-sections.txt
share/examples/rc.d/squid
share/examples/squid/mime.conf
share/examples/squid/squid.conf
@dirrm share/doc/squid
@dirrm share/examples/squid
