# www/squid/PLIST
@comment $NetBSD: PLIST,v 1.2 2001/11/17 23:47:24 kim Exp $
etc/rc.d/squid
share/squid/mime.conf.default
share/squid/squid.conf.default
@unexec rmdir %D/share/squid 2>/dev/null || true
libexec/cachemgr.cgi
sbin/RunAccel
sbin/RunCache
sbin/client
sbin/squid
