# net/gated/PLIST
@comment $NetBSD: PLIST,v 1.3 2005/05/02 20:34:02 reed Exp $
bin/gdc
bin/ospf_monitor
bin/ripquery
man/man8/gated.8
man/man8/gdc.8
man/man8/ospf_monitor.8
man/man8/ripquery.8
sbin/gated
share/doc/gated/config_guide/aggr_stmt.html
share/doc/gated/config_guide/aspath.html
share/doc/gated/config_guide/bgp_stmt.html
share/doc/gated/config_guide/config.html
share/doc/gated/config_guide/control_stmt.html
share/doc/gated/config_guide/defines.html
share/doc/gated/config_guide/direct.html
share/doc/gated/config_guide/dvmrp_stmt.html
share/doc/gated/config_guide/egp_stmt.html
share/doc/gated/config_guide/export_stmt.html
share/doc/gated/config_guide/filter.html
share/doc/gated/config_guide/glossary.html
share/doc/gated/config_guide/icmp_stmt.html
share/doc/gated/config_guide/igmp_stmt.html
share/doc/gated/config_guide/import_stmt.html
share/doc/gated/config_guide/intfcs.html
share/doc/gated/config_guide/intro.html
share/doc/gated/config_guide/isis_stmt.html
share/doc/gated/config_guide/kernel_stmt.html
share/doc/gated/config_guide/options.html
share/doc/gated/config_guide/ospf_stmt.html
share/doc/gated/config_guide/pim_stmt.html
share/doc/gated/config_guide/pref.html
share/doc/gated/config_guide/proto_over.html
share/doc/gated/config_guide/radix.html
share/doc/gated/config_guide/rdisc_stmt.html
share/doc/gated/config_guide/redirect_stmt.html
share/doc/gated/config_guide/references.html
share/doc/gated/config_guide/rip_stmt.html
share/doc/gated/config_guide/samp_eg.html
share/doc/gated/config_guide/samp_host.html
share/doc/gated/config_guide/samp_ig.html
share/doc/gated/config_guide/samp_mcast.html
share/doc/gated/config_guide/snmp_stmt.html
share/doc/gated/config_guide/statements.html
share/doc/gated/config_guide/static_stmt.html
share/doc/gated/config_guide/trace.html
share/doc/gated/copyright-isis.html
share/doc/gated/copyright-ospf.html
share/doc/gated/copyright.html
share/doc/gated/installation_guide/installation.html
share/doc/gated/logo.gif
share/doc/gated/main.html
share/doc/gated/member_info/licenses/advisory.html
share/doc/gated/member_info/licenses/commercial_v4.html
share/doc/gated/member_info/licenses/commercial_v4_new.html
share/doc/gated/member_info/licenses/commercial_v5.html
share/doc/gated/member_info/licenses/commercial_v6.html
share/doc/gated/member_info/licenses/gold.html
share/doc/gated/member_info/licenses/merit.license.v_3.0
share/doc/gated/member_info/licenses/merit_license_v4.html
share/doc/gated/member_info/licenses/research_v4.html
share/doc/gated/member_info/licenses/research_v5.html
share/doc/gated/member_info/licenses/research_v6.html
share/doc/gated/member_info/licenses/service_v4.html
share/doc/gated/member_info/licenses/service_v5.html
share/doc/gated/member_info/licenses/service_v6.html
share/doc/gated/member_info/licenses/small_commercial_v4.html
share/doc/gated/member_info/licenses/small_commercial_v5.html
share/doc/gated/member_info/licenses/small_commercial_v6.html
share/doc/gated/member_info/licenses/small_service_v4.html
share/doc/gated/member_info/licenses/small_service_v5.html
share/doc/gated/member_info/licenses/small_service_v6.html
share/doc/gated/member_info/licenses/supporting_v4.html
share/doc/gated/member_info/licenses/supporting_v5.html
share/doc/gated/member_info/licenses/supporting_v6.html
share/doc/gated/member_info/membership/commercial.member.html
share/doc/gated/member_info/membership/gated-ipv6.summary.html
share/doc/gated/member_info/membership/gated-multi.summary.html
share/doc/gated/member_info/membership/gated-uni.summary.html
share/doc/gated/member_info/membership/goals.html
share/doc/gated/member_info/membership/gold.member.html
share/doc/gated/member_info/membership/membership_categories.html
share/doc/gated/member_info/membership/research.member.html
share/doc/gated/member_info/membership/service.member.html
share/doc/gated/member_info/membership/small_commercial.member.html
share/doc/gated/member_info/membership/small_service.member.html
share/doc/gated/member_info/membership/supporting.member.html
share/doc/gated/member_info/membership/technical_advisory.member.html
share/doc/gated/note.gif
share/doc/gated/operation_guide/gated.html
share/doc/gated/operation_guide/gdc.html
share/doc/gated/operation_guide/operation.html
share/doc/gated/operation_guide/ospf_monitor.html
share/doc/gated/operation_guide/ripquery.html
share/doc/gated/warning.gif
share/examples/rc.d/gated
@dirrm share/doc/gated/config_guide
@dirrm share/doc/gated/installation_guide
@dirrm share/doc/gated/member_info/licenses
@dirrm share/doc/gated/member_info/membership
@dirrm share/doc/gated/member_info
@dirrm share/doc/gated/operation_guide
@dirrm share/doc/gated
