# chat/ircd-hybrid/PLIST
@comment $NetBSD: PLIST,v 1.3 2005/05/02 20:33:58 reed Exp $
man/man8/ircd-hybrid.8
sbin/ircd-hybrid
sbin/ircd-hybrid-fixklines
sbin/ircd-hybrid-mkconf
sbin/ircd-hybrid-mkpasswd
sbin/ircd-hybrid-viconf
sbin/ircd-hybrid-viklines
sbin/ircd-hybrid-vimotd
share/doc/ircd-hybrid/operguide.txt
share/doc/ircd-hybrid/opermyth.txt
share/doc/ircd-hybrid/pkg-setup.txt
share/examples/ircd-hybrid/ircd.conf
share/examples/ircd-hybrid/opers.txt
share/examples/rc.d/ircd-hybrid
@dirrm share/examples/ircd-hybrid
@dirrm share/doc/ircd-hybrid
