# mail/proxsmtp/PLIST
@comment $NetBSD: PLIST,v 1.3 2009/06/14 18:04:40 joerg Exp $
man/man5/proxsmtpd.conf.5
man/man8/proxsmtpd.8
sbin/proxsmtpd
share/examples/proxsmtp/add_header.sh
share/examples/proxsmtp/proxsmtpd.conf
share/examples/proxsmtp/spamassassin.sh
share/examples/rc.d/proxsmtpd
