# net/rabbitmq/files/smf/manifest.xml
<?xml version="1.0"?>
<!DOCTYPE service_bundle SYSTEM "/usr/share/lib/xml/dtd/service_bundle.dtd.1">
<service_bundle type="manifest" name="export">
  <service name="@SMF_PREFIX@/@SMF_NAME@" type="service" version="1">
    <create_default_instance enabled="false" />
    <single_instance />
    <dependency name="network" grouping="require_all" restart_on="error" type="service">
      <service_fmri value="svc:/milestone/network:default" />
    </dependency>
    <dependency name="filesystem" grouping="require_all" restart_on="error" type="service">
      <service_fmri value="svc:/system/filesystem/local" />
    </dependency>
    <dependency name='@SMF_NAME@-epmd' grouping='require_all' restart_on='none' type='service'>
      <service_fmri value='svc:/@SMF_PREFIX@/epmd' />
    </dependency>
    <method_context>
      <method_credential user="@RMQ_USER@" group="@RMQ_GROUP@" />
      <method_environment>
        <envvar name="HOME" value="@RMQ_DIR@" />
	<envvar name="PATH" value="@PREFIX@/sbin:@PREFIX@/bin:/sbin:/usr/sbin:/usr/bin" />
      </method_environment>
    </method_context>
    <exec_method type="method" name="start" exec="@PREFIX@/sbin/rabbitmq-server" timeout_seconds="60" />
    <exec_method type="method" name="stop" exec=":kill" timeout_seconds="60" />
    <property_group name="startd" type="framework">
      <propval name="duration" type="astring" value="child" />
      <propval name="ignore_error" type="astring" value="core,signal" />
    </property_group>
    <template>
      <common_name>
        <loctext xml:lang="C">RabbitMQ Server</loctext>
      </common_name>
    </template>
  </service>
</service_bundle>
