# games/tetrinetx/files/tetrinetx
#!/bin/sh

# PROVIDE: tetrinetx
# REQUIRE: DAEMON LOGIN network

. /etc/rc.subr

name="tetrinetx"
rcvar=$name
command="/usr/pkg/bin/tetrix"
pidfile="/var/run/${name}.pid"
required_files="/usr/pkg/share/tetrinetx/game.conf"

load_rc_config $name
run_rc_command "$1"
