# www/privoxy/PLIST
@comment $NetBSD: PLIST,v 1.3 2005/05/02 20:34:07 reed Exp $
man/man1/privoxy.1
sbin/privoxy
share/doc/privoxy/AUTHORS
share/doc/privoxy/ChangeLog
share/doc/privoxy/LICENSE
share/doc/privoxy/README
share/doc/privoxy/developer-manual/coding.html
share/doc/privoxy/developer-manual/contact.html
share/doc/privoxy/developer-manual/copyright.html
share/doc/privoxy/developer-manual/cvs.html
share/doc/privoxy/developer-manual/documentation.html
share/doc/privoxy/developer-manual/index.html
share/doc/privoxy/developer-manual/introduction.html
share/doc/privoxy/developer-manual/newrelease.html
share/doc/privoxy/developer-manual/quickstart.html
share/doc/privoxy/developer-manual/seealso.html
share/doc/privoxy/developer-manual/testing.html
share/doc/privoxy/developer-manual/webserver-update.html
share/doc/privoxy/faq/configuration.html
share/doc/privoxy/faq/contact.html
share/doc/privoxy/faq/copyright.html
share/doc/privoxy/faq/general.html
share/doc/privoxy/faq/index.html
share/doc/privoxy/faq/installation.html
share/doc/privoxy/faq/misc.html
share/doc/privoxy/faq/trouble.html
share/doc/privoxy/images/files-in-use.jpg
share/doc/privoxy/images/proxy_setup.jpg
share/doc/privoxy/index.html
share/doc/privoxy/man-page/privoxy-man-page.html
share/doc/privoxy/p_doc.css
share/doc/privoxy/user-manual/actions-file.html
share/doc/privoxy/user-manual/appendix.html
share/doc/privoxy/user-manual/config.html
share/doc/privoxy/user-manual/configuration.html
share/doc/privoxy/user-manual/contact.html
share/doc/privoxy/user-manual/copyright.html
share/doc/privoxy/user-manual/filter-file.html
share/doc/privoxy/user-manual/index.html
share/doc/privoxy/user-manual/installation.html
share/doc/privoxy/user-manual/introduction.html
share/doc/privoxy/user-manual/quickstart.html
share/doc/privoxy/user-manual/seealso.html
share/doc/privoxy/user-manual/startup.html
share/doc/privoxy/user-manual/templates.html
share/doc/privoxy/user-manual/upgradersnote.html
share/examples/privoxy/config
share/examples/privoxy/default.action
share/examples/privoxy/default.filter
share/examples/privoxy/standard.action
share/examples/privoxy/templates/blocked
share/examples/privoxy/templates/cgi-error-404
share/examples/privoxy/templates/cgi-error-bad-param
share/examples/privoxy/templates/cgi-error-disabled
share/examples/privoxy/templates/cgi-error-file
share/examples/privoxy/templates/cgi-error-file-read-only
share/examples/privoxy/templates/cgi-error-modified
share/examples/privoxy/templates/cgi-error-parse
share/examples/privoxy/templates/cgi-style.css
share/examples/privoxy/templates/connect-failed
share/examples/privoxy/templates/default
share/examples/privoxy/templates/edit-actions-add-url-form
share/examples/privoxy/templates/edit-actions-for-url
share/examples/privoxy/templates/edit-actions-for-url-filter
share/examples/privoxy/templates/edit-actions-list
share/examples/privoxy/templates/edit-actions-list-button
share/examples/privoxy/templates/edit-actions-list-section
share/examples/privoxy/templates/edit-actions-list-url
share/examples/privoxy/templates/edit-actions-remove-url-form
share/examples/privoxy/templates/edit-actions-url-form
share/examples/privoxy/templates/mod-local-help
share/examples/privoxy/templates/mod-support-and-service
share/examples/privoxy/templates/mod-title
share/examples/privoxy/templates/mod-unstable-warning
share/examples/privoxy/templates/no-such-domain
share/examples/privoxy/templates/show-request
share/examples/privoxy/templates/show-status
share/examples/privoxy/templates/show-status-file
share/examples/privoxy/templates/show-url-info
share/examples/privoxy/templates/show-version
share/examples/privoxy/templates/toggle
share/examples/privoxy/templates/toggle-mini
share/examples/privoxy/templates/untrusted
share/examples/privoxy/trust
share/examples/privoxy/user.action
share/examples/rc.d/privoxy
@dirrm share/examples/privoxy/templates
@dirrm share/examples/privoxy
@dirrm share/doc/privoxy/user-manual
@dirrm share/doc/privoxy/man-page
@dirrm share/doc/privoxy/images
@dirrm share/doc/privoxy/faq
@dirrm share/doc/privoxy/developer-manual
@dirrm share/doc/privoxy
