# net/powerdns/PLIST
@comment $NetBSD: PLIST,v 1.2 2009/06/14 18:09:38 joerg Exp $
bin/pdns_control
bin/zone2ldap
bin/zone2sql
lib/pdns/libgeobackend.la
lib/pdns/libpipebackend.la
man/man8/pdns_control.8
man/man8/pdns_server.8
man/man8/zone2sql.8
sbin/pdns_server
share/examples/pdns/pdns.conf-dist
share/examples/rc.d/pdns
