# japanese/freewnn-lib/files/ja_freewnn
#!/bin/sh
#
# $NetBSD: ja_freewnn,v 1.1 2002/02/28 17:45:45 wiz Exp $
#
# PROVIDE: freewnn
# REQUIRE: DAEMON

. /etc/rc.subr

name="freewnn"
command_start="@PREFIX@/sbin/jserver"
command_stop="@PREFIX@/sbin/wnnkill"
pidfile=""

cmd=${1:-start}

case ${cmd} in
start)
    if [ -x ${command_start} ]; then
	echo "Starting ${name}."
        ${command_start} >/dev/null 2>&1
    fi
    ;;
stop)
    if [ -x ${command_stop} ]; then
        echo "Stopping ${name}."
        ${command_stop} >/dev/null 2>&1
    fi
    ;;
*)
    :
    ;;
esac
exit 0
