# security/prelude-pflogger/PLIST
@comment $NetBSD: PLIST,v 1.2 2008/07/21 12:22:15 shannonjr Exp $
bin/prelude-pflogger
sbin/run-prelude-pflogger
share/examples/rc.d/pflogger
share/examples/prelude-pflogger/prelude-pflogger.conf
