# www/wwwoffle/PLIST
@comment $NetBSD: PLIST,v 1.1 2001/11/01 02:19:10 zuntum Exp $
bin/wwwoffle
bin/wwwoffle-ls
bin/wwwoffle-mv
bin/wwwoffle-read
bin/wwwoffle-rm
bin/wwwoffle-tools
bin/wwwoffle-write
etc/rc.d/wwwoffle
etc/rc.d/wwwoffled
etc/wwwoffle/wwwoffle.conf.install
@exec if [ ! -f %D/etc/wwwoffle/wwwoffle.conf ]; then cp -p %D/etc/wwwoffle/wwwoffle.conf.install %D/etc/wwwoffle/wwwoffle.conf; fi
libdata/wwwoffle-spool.tar.gz
man/man1/wwwoffle.1
man/man5/wwwoffle.conf.5
man/man8/wwwoffled.8
sbin/wwwoffle-convert-cache
sbin/wwwoffle-upgrade-config.pl
sbin/wwwoffle-upgrade-config-2.x-2.5.pl
sbin/wwwoffled
share/doc/wwwoffle/CHANGES.CONF
share/doc/wwwoffle/CONVERT
share/doc/wwwoffle/COPYING
share/doc/wwwoffle/FAQ
share/doc/wwwoffle/INSTALL
share/doc/wwwoffle/NEWS
share/doc/wwwoffle/README
share/doc/wwwoffle/README.1st
share/doc/wwwoffle/README.CONF
share/doc/wwwoffle/README.PWD
share/doc/wwwoffle/README.URL
share/doc/wwwoffle/README.compress
share/doc/wwwoffle/README.htdig
share/doc/wwwoffle/README.lang
share/doc/wwwoffle/README.udmsearch
share/doc/wwwoffle/README.win32
@dirrm share/doc/wwwoffle
@exec ${GTAR} xzCf / %D/libdata/wwwoffle-spool.tar.gz
@exec cd /var/wwwoffle/html && mv wwwoffle.pac wwwoffle.pac.dist && sed /PROXY/s/localhost/`hostname`/ <wwwoffle.pac.dist >wwwoffle.pac
@exec ${MKDIR} /var/wwwoffle/outgoing
@exec ${MKDIR} /var/wwwoffle/monitor
@exec ${MKDIR} /var/wwwoffle/lastime
@exec ${MKDIR} /var/wwwoffle/prevtime1
@exec ${PREFIX}/sbin/wwwoffle-convert-cache /var/wwwoffle
@unexec cd /var/wwwoffle && rm -rf html*
@unexec rmdir /var/wwwoffle/* 2>/dev/null && rmdir /var/wwwoffle 2>/dev/null || true
@unexec if [ -d /var/wwwoffle ]; then echo "If you won't use the WWWOFFLE package anymore, remove the cache manually by"; echo '	rm -rf /var/wwwoffle'; fi
