# net/ruby-stompserver/PLIST
@comment $NetBSD: PLIST,v 1.2 2012/08/13 18:43:33 imil Exp $
bin/stompserver
${GEM_HOME}/cache/${GEM_NAME}.gem
${GEM_LIBDIR}/.gemtest
${GEM_LIBDIR}/History.txt
${GEM_LIBDIR}/Manifest.txt
${GEM_LIBDIR}/README.txt
${GEM_LIBDIR}/Rakefile
${GEM_LIBDIR}/STATUS
${GEM_LIBDIR}/bin/stompserver
${GEM_LIBDIR}/client/README.txt
${GEM_LIBDIR}/client/both.rb
${GEM_LIBDIR}/client/consume.rb
${GEM_LIBDIR}/client/send.rb
${GEM_LIBDIR}/config/stompserver.conf
${GEM_LIBDIR}/etc/passwd.example
${GEM_LIBDIR}/lib/stomp_server.rb
${GEM_LIBDIR}/lib/stomp_server/protocols/http.rb
${GEM_LIBDIR}/lib/stomp_server/protocols/stomp.rb
${GEM_LIBDIR}/lib/stomp_server/queue.rb
${GEM_LIBDIR}/lib/stomp_server/queue/activerecord_queue.rb
${GEM_LIBDIR}/lib/stomp_server/queue/ar_message.rb
${GEM_LIBDIR}/lib/stomp_server/queue/dbm_queue.rb
${GEM_LIBDIR}/lib/stomp_server/queue/file_queue.rb
${GEM_LIBDIR}/lib/stomp_server/queue/memory_queue.rb
${GEM_LIBDIR}/lib/stomp_server/queue_manager.rb
${GEM_LIBDIR}/lib/stomp_server/stomp_auth.rb
${GEM_LIBDIR}/lib/stomp_server/stomp_frame.rb
${GEM_LIBDIR}/lib/stomp_server/stomp_id.rb
${GEM_LIBDIR}/lib/stomp_server/stomp_user.rb
${GEM_LIBDIR}/lib/stomp_server/test_server.rb
${GEM_LIBDIR}/lib/stomp_server/topic_manager.rb
${GEM_LIBDIR}/setup.rb
${GEM_LIBDIR}/test/tesly.rb
${GEM_LIBDIR}/test/test_queue_manager.rb
${GEM_LIBDIR}/test/test_stomp_frame.rb
${GEM_LIBDIR}/test/test_topic_manager.rb
${GEM_LIBDIR}/test_todo/test_stomp_server.rb
${GEM_HOME}/specifications/${GEM_NAME}.gemspec
share/examples/rc.d/stompserver
share/examples/stompserver/stompserver.conf
