# net/xymonclient/PLIST
@comment $NetBSD: PLIST,v 1.5 2012/08/15 16:04:44 spz Exp $
libexec/xymon/client/clientupdate
libexec/xymon/client/local/README
libexec/xymon/client/logfetch
libexec/xymon/client/msgcache
${PLIST.meminfo}libexec/xymon/client/${LOWER_OPSYS}-meminfo
libexec/xymon/client/orcaxymon
libexec/xymon/client/runclient.sh
libexec/xymon/client/xymon
libexec/xymon/client/xymoncfg
libexec/xymon/client/xymonclient-aix.sh
libexec/xymon/client/xymonclient-darwin.sh
libexec/xymon/client/xymonclient-dragonfly.sh
libexec/xymon/client/xymonclient-freebsd.sh
libexec/xymon/client/xymonclient-hp-ux.sh
libexec/xymon/client/xymonclient-irix.sh
libexec/xymon/client/xymonclient-linux.sh
libexec/xymon/client/xymonclient-netbsd.sh
libexec/xymon/client/xymonclient-openbsd.sh
libexec/xymon/client/xymonclient-osf1.sh
libexec/xymon/client/xymonclient-sco_sv.sh
libexec/xymon/client/xymonclient-sunos.sh
libexec/xymon/client/xymonclient-unixware.sh
libexec/xymon/client/xymonclient.sh
libexec/xymon/client/xymoncmd
libexec/xymon/client/xymondigest
libexec/xymon/client/xymongrep
libexec/xymon/client/xymonlaunch
share/examples/rc.d/xymonclient
share/examples/xymon/clientlaunch.cfg
share/examples/xymon/localclient.cfg
share/examples/xymon/xymonclient.cfg
@pkgdir share/examples/xymon
@pkgdir libexec/xymon/ext/client
