# net/xymonclient/PLIST
@comment $NetBSD: PLIST,v 1.1.1.1 2010/02/14 16:11:11 spz Exp $
libexec/xymon/client/bb
libexec/xymon/client/bbcmd
libexec/xymon/client/bbdigest
libexec/xymon/client/bbhostgrep
libexec/xymon/client/bbhostshow
libexec/xymon/client/clientupdate
libexec/xymon/client/hobbitclient-aix.sh
libexec/xymon/client/hobbitclient-darwin.sh
libexec/xymon/client/hobbitclient-freebsd.sh
libexec/xymon/client/hobbitclient-hp-ux.sh
libexec/xymon/client/hobbitclient-irix.sh
libexec/xymon/client/hobbitclient-linux.sh
libexec/xymon/client/hobbitclient-${LOWER_OPSYS}.sh
libexec/xymon/client/hobbitclient-openbsd.sh
libexec/xymon/client/hobbitclient-osf1.sh
libexec/xymon/client/hobbitclient-sco_sv.sh
libexec/xymon/client/hobbitclient-sunos.sh
libexec/xymon/client/hobbitclient.sh
libexec/xymon/client/hobbitlaunch
libexec/xymon/client/logfetch
libexec/xymon/client/msgcache
libexec/xymon/client/${LOWER_OPSYS}-meminfo
libexec/xymon/client/orcahobbit
libexec/xymon/client/runclient.sh
share/examples/rc.d/xymonclient
share/examples/xymon/clientlaunch.cfg
share/examples/xymon/hobbitclient.cfg
share/examples/xymon/localclient.cfg
