# security/honeyd-arpd/PLIST
@comment $NetBSD: PLIST,v 1.3 2005/05/02 20:34:05 reed Exp $
man/man8/honeyd-arpd.8
sbin/honeyd-arpd
share/examples/rc.d/honeydarpd
