# security/winbind/PLIST
@comment $NetBSD: PLIST,v 1.1.1.1 2002/07/29 04:19:00 jlam Exp $
bin/wbinfo
etc/rc.d/winbindd
lib/libnss_winbind.so
lib/security/pam_winbind.so
sbin/winbindd
