# mail/enma/PLIST
@comment $NetBSD: PLIST,v 1.2 2009/04/09 02:25:37 obache Exp $
bin/sidfquery
libexec/enma
man/ja_JP.UTF-8/man1/enma.1
man/man1/enma.1
share/examples/enma/enma.conf.sample
share/examples/rc.d/enma
@dirrm share/examples/enma
