# net/nagios-nrpe/PLIST
@comment $NetBSD: PLIST,v 1.3 2006/05/26 18:53:25 wiz Exp $
sbin/nrpe
libexec/nagios/check_nrpe
share/examples/nagios/nrpe.cfg
share/examples/rc.d/nrpe
@unexec ${RMDIR} %D/share/examples/nagios 2>/dev/null || ${TRUE}
@unexec ${RMDIR} %D/libexec/nagios 2>/dev/null || ${TRUE}
