# net/nagios-nrpe/PLIST
@comment $NetBSD: PLIST,v 1.4 2009/06/14 21:00:04 joerg Exp $
sbin/nrpe
libexec/nagios/check_nrpe
share/examples/nagios/nrpe.cfg
share/examples/rc.d/nrpe
