# sysutils/syslog-ng/PLIST
@comment $NetBSD: PLIST,v 1.3 2011/09/21 19:33:09 bouyer Exp $

bin/loggen
bin/pdbtool
bin/syslog2ng
bin/update-patterndb
etc/modules.conf
etc/scl.conf
etc/syslog-ng.conf
lib/libsyslog-ng.la
lib/libsyslog-ng.so
lib/libsyslog-ng.so.0
lib/libsyslog-ng.so.0.0
lib/syslog-ng/libaffile.la
lib/syslog-ng/libaffile.so
lib/syslog-ng/libafprog.la
lib/syslog-ng/libafprog.so
lib/syslog-ng/libafsocket-notls.la
lib/syslog-ng/libafsocket-notls.so
lib/syslog-ng/libafsocket-tls.la
lib/syslog-ng/libafsocket-tls.so
lib/syslog-ng/libafsocket.so
lib/syslog-ng/libafuser.la
lib/syslog-ng/libafuser.so
lib/syslog-ng/libbasicfuncs.la
lib/syslog-ng/libbasicfuncs.so
lib/syslog-ng/libconfgen.la
lib/syslog-ng/libconfgen.so
lib/syslog-ng/libconvertfuncs.la
lib/syslog-ng/libconvertfuncs.so
lib/syslog-ng/libcsvparser.la
lib/syslog-ng/libcsvparser.so
lib/syslog-ng/libdbparser.la
lib/syslog-ng/libdbparser.so
lib/syslog-ng/libdummy.la
lib/syslog-ng/libdummy.so
lib/syslog-ng/libsyslogformat.la
lib/syslog-ng/libsyslogformat.so
man/man1/loggen.1
man/man1/pdbtool.1
man/man1/syslog-ng-ctl.1
man/man5/syslog-ng.conf.5
man/man8/syslog-ng.8
sbin/syslog-ng
sbin/syslog-ng-ctl
share/doc/${PKGBASE}/INSTALL
share/doc/${PKGBASE}/README
share/doc/${PKGBASE}/syslog-ng.conf.doc
share/examples/rc.d/syslogng
share/examples/syslog-ng/syslog-ng.conf.sample
share/examples/syslog-ng/syslog-ng.conf.solaris
share/examples/syslog-ng/syslog-ng.service
share/include/scl/pacct/plugin.conf
share/include/scl/syslogconf/README
share/include/scl/syslogconf/convert-syslogconf.awk
share/include/scl/syslogconf/plugin.conf
share/include/scl/system/generate-system-source.sh
share/include/scl/system/plugin.conf
share/xsd/patterndb-1.xsd
share/xsd/patterndb-2.xsd
share/xsd/patterndb-3.xsd
