# sysutils/syslog-ng/PLIST
@comment $NetBSD: PLIST,v 1.2 2009/06/14 18:16:14 joerg Exp $
bin/syslog2ng
man/man5/syslog-ng.conf.5
man/man8/syslog-ng.8
sbin/syslog-ng
share/doc/${PKGBASE}/book1.html
share/doc/${PKGBASE}/destinations.html
share/doc/${PKGBASE}/intro.html
share/doc/${PKGBASE}/logpath.html
share/doc/${PKGBASE}/msgroute.html
share/doc/${PKGBASE}/reference.html
share/doc/${PKGBASE}/tuning.html
share/doc/${PKGBASE}/index.html
share/doc/${PKGBASE}/x731.html
share/doc/${PKGBASE}/x771.html
share/doc/${PKGBASE}/x903.html
share/doc/${PKGBASE}/x909.html
share/doc/${PKGBASE}/x172.html
share/doc/${PKGBASE}/x361.html
share/doc/${PKGBASE}/x97.html
share/doc/${PKGBASE}/INSTALL
share/doc/${PKGBASE}/README
share/doc/${PKGBASE}/README.spoof
share/doc/${PKGBASE}/syslog-ng.conf.doc
share/examples/rc.d/syslogng
share/examples/syslog-ng/syslog-ng.conf.demo
share/examples/syslog-ng/syslog-ng.conf.sample
share/examples/syslog-ng/syslog-ng.conf.solaris
