# security/honeyd/PLIST
@comment $NetBSD: PLIST,v 1.5 2005/05/02 20:34:04 reed Exp $
bin/honeyd
bin/honeydctl
include/honeyd/debug.h
include/honeyd/hooks.h
include/honeyd/plugins.h
include/honeyd/plugins_config.h
lib/honeyd/libhoneyd.so
man/man1/honeydctl.1
man/man8/honeyd.8
share/doc/honeyd/README
share/examples/honeyd/config.ethernet
share/examples/honeyd/config.sample
share/examples/honeyd/nmap.assoc
share/examples/honeyd/nmap.prints
share/examples/honeyd/pf.os
share/examples/honeyd/webserver/htdocs/config.py
share/examples/honeyd/webserver/htdocs/graphs/traffic_daily.gif
share/examples/honeyd/webserver/htdocs/graphs/traffic_hourly.gif
share/examples/honeyd/webserver/htdocs/images/add.gif
share/examples/honeyd/webserver/htdocs/images/arrow_down.gif
share/examples/honeyd/webserver/htdocs/images/arrow_right.gif
share/examples/honeyd/webserver/htdocs/images/delete.gif
share/examples/honeyd/webserver/htdocs/images/edit.gif
share/examples/honeyd/webserver/htdocs/images/google_logo.gif
share/examples/honeyd/webserver/htdocs/images/logo.gif
share/examples/honeyd/webserver/htdocs/index.py
share/examples/honeyd/webserver/htdocs/styles/layout.css
share/examples/honeyd/webserver/htdocs/templates/config_ip.tmpl
share/examples/honeyd/webserver/htdocs/templates/inc/footer.tmpl
share/examples/honeyd/webserver/htdocs/templates/inc/header.tmpl
share/examples/honeyd/webserver/htdocs/templates/inc/main.tmpl
share/examples/honeyd/webserver/htdocs/templates/inc/nav.tmpl
share/examples/honeyd/webserver/htdocs/templates/index.tmpl
share/examples/honeyd/webserver/htdocs/templates/status_connections.tmpl
share/examples/honeyd/webserver/htdocs/templates/status_stats.tmpl
share/examples/honeyd/webserver/htmltmpl.py
share/examples/honeyd/webserver/server.py
share/examples/honeyd/webserver/support.py
share/examples/honeyd/xprobe2.conf
share/examples/rc.d/honeyd
@dirrm share/examples/honeyd/webserver/htdocs/templates/inc
@dirrm share/examples/honeyd/webserver/htdocs/templates
@dirrm share/examples/honeyd/webserver/htdocs/styles
@dirrm share/examples/honeyd/webserver/htdocs/images
@dirrm share/examples/honeyd/webserver/htdocs/graphs
@dirrm share/examples/honeyd/webserver/htdocs
@dirrm share/examples/honeyd/webserver
@dirrm share/examples/honeyd
@dirrm share/doc/honeyd
@dirrm lib/honeyd
@dirrm include/honeyd
