# www/instiki/PLIST
@comment $NetBSD: PLIST,v 1.2 2005/05/02 20:34:07 reed Exp $
share/doc/instiki/README
share/examples/rc.d/instiki
share/instiki/app/controllers/wiki.rb
share/instiki/app/models/author.rb
share/instiki/app/models/chunks/category.rb
share/instiki/app/models/chunks/category_test.rb
share/instiki/app/models/chunks/chunk.rb
share/instiki/app/models/chunks/engines.rb
share/instiki/app/models/chunks/include.rb
share/instiki/app/models/chunks/literal.rb
share/instiki/app/models/chunks/match.rb
share/instiki/app/models/chunks/nowiki.rb
share/instiki/app/models/chunks/nowiki_test.rb
share/instiki/app/models/chunks/test.rb
share/instiki/app/models/chunks/uri.rb
share/instiki/app/models/chunks/uri_test.rb
share/instiki/app/models/chunks/wiki.rb
share/instiki/app/models/chunks/wiki_test.rb
share/instiki/app/models/page.rb
share/instiki/app/models/page_lock.rb
share/instiki/app/models/page_set.rb
share/instiki/app/models/page_test.rb
share/instiki/app/models/revision.rb
share/instiki/app/models/revision_test.rb
share/instiki/app/models/web.rb
share/instiki/app/models/web_test.rb
share/instiki/app/models/wiki_content.rb
share/instiki/app/models/wiki_service.rb
share/instiki/app/models/wiki_service_test.rb
share/instiki/app/models/wiki_words.rb
share/instiki/app/models/wiki_words_test.rb
share/instiki/app/views/bottom.rhtml
share/instiki/app/views/markdown_help.rhtml
share/instiki/app/views/navigation.rhtml
share/instiki/app/views/rdoc_help.rhtml
share/instiki/app/views/static_style_sheet.rhtml
share/instiki/app/views/textile_help.rhtml
share/instiki/app/views/top.rhtml
share/instiki/app/views/wiki/authors.rhtml
share/instiki/app/views/wiki/edit.rhtml
share/instiki/app/views/wiki/edit_web.rhtml
share/instiki/app/views/wiki/export.rhtml
share/instiki/app/views/wiki/feeds.rhtml
share/instiki/app/views/wiki/list.rhtml
share/instiki/app/views/wiki/locked.rhtml
share/instiki/app/views/wiki/login.rhtml
share/instiki/app/views/wiki/new.rhtml
share/instiki/app/views/wiki/new_system.rhtml
share/instiki/app/views/wiki/new_web.rhtml
share/instiki/app/views/wiki/page.rhtml
share/instiki/app/views/wiki/print.rhtml
share/instiki/app/views/wiki/published.rhtml
share/instiki/app/views/wiki/recently_revised.rhtml
share/instiki/app/views/wiki/revision.rhtml
share/instiki/app/views/wiki/rollback.rhtml
share/instiki/app/views/wiki/rss_feed.rhtml
share/instiki/app/views/wiki/search.rhtml
share/instiki/app/views/wiki/tex.rhtml
share/instiki/app/views/wiki/tex_web.rhtml
share/instiki/app/views/wiki/web_list.rhtml
share/instiki/app/views/wiki_words_help.rhtml
share/instiki/instiki
share/instiki/libraries/action_controller_servlet.rb
share/instiki/libraries/bluecloth.rb
share/instiki/libraries/diff/diff.rb
share/instiki/libraries/diff/diff_test.rb
share/instiki/libraries/erb.rb
share/instiki/libraries/madeleine.rb
share/instiki/libraries/madeleine/automatic.rb
share/instiki/libraries/madeleine/clock.rb
share/instiki/libraries/madeleine/files.rb
share/instiki/libraries/madeleine/zmarshal.rb
share/instiki/libraries/madeleine_service.rb
share/instiki/libraries/rdocsupport.rb
share/instiki/libraries/redcloth.rb
share/instiki/libraries/redcloth_for_tex.rb
share/instiki/libraries/redcloth_for_tex_test.rb
share/instiki/libraries/view_helper.rb
share/instiki/libraries/web_controller_server.rb
share/instiki/libraries/zip/ioextras.rb
share/instiki/libraries/zip/stdrubyext.rb
share/instiki/libraries/zip/tempfile_bugfixed.rb
share/instiki/libraries/zip/zip.rb
share/instiki/libraries/zip/zipfilesystem.rb
share/instiki/libraries/zip/ziprequire.rb
@dirrm share/instiki/libraries/zip
@dirrm share/instiki/libraries/madeleine
@dirrm share/instiki/libraries/diff
@dirrm share/instiki/libraries
@dirrm share/instiki/app/views/wiki
@dirrm share/instiki/app/views
@dirrm share/instiki/app/models/chunks
@dirrm share/instiki/app/models
@dirrm share/instiki/app/controllers
@dirrm share/instiki/app
@dirrm share/instiki
@dirrm share/doc/instiki
