# net/isc-dhcpd/PLIST
@comment $NetBSD: PLIST,v 1.3 2009/06/14 18:09:31 joerg Exp $
bin/omshell
include/dhcpctl.h
lib/libdhcpctl.a
lib/libomapi.a
man/man1/omshell.1
man/man3/dhcpctl.3
man/man3/omapi.3
man/man3/omshell.3
man/man5/dhcpd.conf.5
man/man5/dhcpd.leases.5
man/man8/dhcpd.8
sbin/dhcpd
share/examples/isc-dhcp-server/dhcpd.conf
share/examples/rc.d/isc_dhcpd
