# mail/courier-imap/PLIST
@comment $NetBSD: PLIST,v 1.12 2005/05/02 20:33:59 reed Exp $
bin/couriertls
bin/imapd
libexec/courier/courierpop3d
libexec/courier/courierpop3login
libexec/courier/imaplogin
man/man1/couriertcpd.1
man/man8/courierimapd.8
man/man8/courierpop3d.8
man/man8/mkimapdcert.8
man/man8/mkpop3dcert.8
sbin/couriertcpd
sbin/imapd
sbin/imapd-ssl
sbin/mkimapdcert
sbin/mkpop3dcert
sbin/pop3d
sbin/pop3d-ssl
share/courier/imapd
share/courier/imapd-ssl
share/courier/mkimapdcert
share/courier/mkpop3dcert
share/courier/pop3d
share/courier/pop3d-ssl
share/doc/courier/imap/BUGS.html
share/doc/courier/imap/README.html
share/doc/courier/imap/README.proxy.html
share/examples/courier/imap.pam
share/examples/courier/imapd-ssl.dist
share/examples/courier/imapd.cnf
share/examples/courier/imapd.dist
share/examples/courier/pop3.pam
share/examples/courier/pop3d-ssl.dist
share/examples/courier/pop3d.cnf
share/examples/courier/pop3d.dist
share/examples/rc.d/courierpops
share/examples/rc.d/courierpop
share/examples/rc.d/courierimaps
share/examples/rc.d/courierimap
@unexec ${RMDIR} -p %D/share/examples/courier 2>/dev/null || ${TRUE}
@dirrm share/doc/courier/imap
@unexec ${RMDIR} -p %D/share/doc/courier 2>/dev/null || ${TRUE}
@unexec ${RMDIR} -p %D/share/courier 2>/dev/null || ${TRUE}
@unexec ${RMDIR} -p %D/libexec/courier 2>/dev/null || ${TRUE}
