# mail/courier-imap/PLIST
@comment $NetBSD: PLIST,v 1.1 2001/11/01 00:59:42 zuntum Exp $
bin/couriertls
bin/imapd
bin/imapd.maildirmake
bin/pop3d
etc/rc.d/courierimap
etc/rc.d/courierimaps
etc/rc.d/courierpop
etc/rc.d/courierpops
libexec/courier/authlib/authcram
libexec/courier/authlib/authcustom
libexec/courier/authlib/authpwd
libexec/courier/authlib/authuserdb
libexec/courier/couriertcpd
libexec/courier/deliverquota
libexec/courier/imapd-ssl.rc
libexec/courier/imapd.rc
libexec/courier/logger
libexec/courier/makedatprog
libexec/courier/pop3d-ssl.rc
libexec/courier/pop3d.rc
man/man1/imapd.maildirmake.1
man/man8/authcram.8
man/man8/authdaemon.8
man/man8/authdaemond.8
man/man8/authldap.8
man/man8/authlib.8
man/man8/authmysql.8
man/man8/authpam.8
man/man8/authpwd.8
man/man8/authshadow.8
man/man8/authuserdb.8
man/man8/authvchkpw.8
man/man8/courier-imapd.8
man/man8/couriertcpd.8
man/man8/deliverquota.8
man/man8/makeuserdb.8
man/man8/mkimapdcert.8
man/man8/mkpop3dcert.8
man/man8/pw2userdb.8
man/man8/userdb.8
man/man8/userdbpw.8
man/man8/vchkpw2userdb.8
sbin/imaplogin
sbin/makeuserdb
sbin/mkimapdcert
sbin/mkpop3dcert
sbin/pop3login
sbin/pw2userdb
sbin/sysconftool
sbin/userdb
sbin/userdbpw
sbin/vchkpw2userdb
share/courier/makeuserdb
share/courier/mkimapdcert
share/courier/mkpop3dcert
share/courier/pw2userdb
share/courier/userdb
share/courier/vchkpw2userdb
share/doc/courier/README.imap
share/doc/courier/README.maildirquota
share/doc/courier/README.sharedfolders
share/examples/courier/imapd-ssl.dist
share/examples/courier/imapd.cnf
share/examples/courier/imapd.dist
share/examples/courier/pop3d-ssl.dist
share/examples/courier/pop3d.cnf
share/examples/courier/pop3d.dist
share/examples/courier/quotawarnmsg.example
@unexec rmdir -p %D/share/examples/courier 2>/dev/null || true
@unexec rmdir -p %D/share/doc/courier 2>/dev/null || true
@unexec rmdir -p %D/share/courier 2>/dev/null || true
@unexec rmdir -p %D/libexec/courier 2>/dev/null || true
