# mail/courier-imap/PLIST
@comment $NetBSD: PLIST,v 1.16 2009/06/14 18:04:30 joerg Exp $
bin/imapd
libexec/courier/courierpop3d
libexec/courier/courierpop3login
libexec/courier/imaplogin
man/man8/courierimapd.8
man/man8/courierpop3d.8
man/man8/mkimapdcert.8
man/man8/mkpop3dcert.8
sbin/imapd
sbin/imapd-ssl
sbin/mkimapdcert
sbin/mkpop3dcert
sbin/pop3d
sbin/pop3d-ssl
share/courier/imapd
share/courier/imapd-ssl
share/courier/mkimapdcert
share/courier/mkpop3dcert
share/courier/pop3d
share/courier/pop3d-ssl
share/doc/courier/courierpop3d.html
share/doc/courier/imap/BUGS.html
share/doc/courier/imap/README.html
share/doc/courier/imap/README.proxy.html
share/doc/courier/imapd.html
share/doc/courier/mkimapdcert.html
share/doc/courier/mkpop3dcert.html
share/examples/courier/imap.authpam
share/examples/courier/imapd-ssl.dist
share/examples/courier/imapd.cnf
share/examples/courier/imapd.dist
share/examples/courier/pop3.authpam
share/examples/courier/pop3d-ssl.dist
share/examples/courier/pop3d.cnf
share/examples/courier/pop3d.dist
share/examples/rc.d/courierimap
share/examples/rc.d/courierimaps
share/examples/rc.d/courierpop
share/examples/rc.d/courierpops
