# net/snort/PLIST
@comment $NetBSD: PLIST,v 1.7.2.1 2003/03/07 07:46:37 jmc Exp $
bin/snort
etc/rc.d/snort
man/man8/snort.8
share/doc/snort/AUTHORS
share/doc/snort/BUGS
share/doc/snort/CREDITS
share/doc/snort/FAQ
share/doc/snort/INSTALL
share/doc/snort/NEWS
share/doc/snort/README
share/doc/snort/README.FLEXRESP
share/doc/snort/README.PLUGINS
share/doc/snort/README.SNMP
share/doc/snort/README.csv
share/doc/snort/README.database
share/doc/snort/README.xml
share/doc/snort/RULES.todo
share/doc/snort/SnortUsersManual.pdf
share/doc/snort/TODO
share/doc/snort/USAGE
share/examples/snort/classification.config
share/examples/snort/reference.config
share/examples/snort/snort.conf.default
share/snort/rules/SnortCommonMIB.txt
share/snort/rules/SnortIDAlertMIB.txt
share/snort/rules/attack-responses.rules
share/snort/rules/backdoor.rules
share/snort/rules/bad-traffic.rules
share/snort/rules/chat.rules
share/snort/rules/ddos.rules
share/snort/rules/deleted.rules
share/snort/rules/dns.rules
share/snort/rules/dos.rules
share/snort/rules/experimental.rules
share/snort/rules/exploit.rules
share/snort/rules/finger.rules
share/snort/rules/ftp.rules
share/snort/rules/gen-msg.map
share/snort/rules/icmp-info.rules
share/snort/rules/icmp.rules
share/snort/rules/imap.rules
share/snort/rules/info.rules
share/snort/rules/local.rules
share/snort/rules/misc.rules
share/snort/rules/multimedia.rules
share/snort/rules/mysql.rules
share/snort/rules/netbios.rules
share/snort/rules/nntp.rules
share/snort/rules/oracle.rules
share/snort/rules/other-ids.rules
share/snort/rules/p2p.rules
share/snort/rules/policy.rules
share/snort/rules/pop2.rules
share/snort/rules/pop3.rules
share/snort/rules/porn.rules
share/snort/rules/rpc.rules
share/snort/rules/rservices.rules
share/snort/rules/scan.rules
share/snort/rules/shellcode.rules
share/snort/rules/sid
share/snort/rules/sid-msg.map
share/snort/rules/smtp.rules
share/snort/rules/snmp.rules
share/snort/rules/sql.rules
share/snort/rules/telnet.rules
share/snort/rules/tftp.rules
share/snort/rules/virus.rules
share/snort/rules/web-attacks.rules
share/snort/rules/web-cgi.rules
share/snort/rules/web-client.rules
share/snort/rules/web-coldfusion.rules
share/snort/rules/web-frontpage.rules
share/snort/rules/web-iis.rules
share/snort/rules/web-misc.rules
share/snort/rules/web-php.rules
share/snort/rules/x11.rules
@dirrm share/snort/rules
@dirrm share/examples/snort
@dirrm share/doc/snort
@dirrm share/snort
