# net/snort/PLIST
@comment $NetBSD: PLIST,v 1.26 2008/01/06 00:28:44 adrianp Exp $
bin/snort
lib/snort_dynamicengine/libsf_engine.la
lib/snort_dynamicrules/lib_sfdynamic_example_rule.la
lib/snort_dynamicpreprocessor/libsf_dcerpc_preproc.la
lib/snort_dynamicpreprocessor/libsf_dns_preproc.la
lib/snort_dynamicpreprocessor/libsf_ftptelnet_preproc.la
lib/snort_dynamicpreprocessor/libsf_smtp_preproc.la
lib/snort_dynamicpreprocessor/libsf_ssh_preproc.la
lib/snort_dynamicpreprocessor/lib_sfdynamic_preprocessor_example.la
man/man8/snort.8
share/doc/snort/AUTHORS
share/doc/snort/BUGS
share/doc/snort/CREDITS
share/doc/snort/INSTALL
share/doc/snort/NEWS
share/doc/snort/PROBLEMS
share/doc/snort/README
share/doc/snort/README.ARUBA
share/doc/snort/README.FLEXRESP
share/doc/snort/README.FLEXRESP2
share/doc/snort/README.INLINE
share/doc/snort/README.PLUGINS
share/doc/snort/README.PerfProfiling
share/doc/snort/README.SMTP
share/doc/snort/README.UNSOCK
share/doc/snort/README.WIN32
share/doc/snort/README.alert_order
share/doc/snort/README.asn1
share/doc/snort/README.csv
share/doc/snort/README.database
share/doc/snort/README.dcerpc
share/doc/snort/README.decode
share/doc/snort/README.dns
share/doc/snort/README.event_queue
share/doc/snort/README.flow
share/doc/snort/README.flow-portscan
share/doc/snort/README.flowbits
share/doc/snort/README.frag3
share/doc/snort/README.ftptelnet
share/doc/snort/README.http_inspect
share/doc/snort/README.ipv6
share/doc/snort/README.ppm
share/doc/snort/README.sfportscan
share/doc/snort/README.ssh
share/doc/snort/README.stream4
share/doc/snort/README.stream5
share/doc/snort/README.tag
share/doc/snort/README.thresholding
share/doc/snort/README.variables
share/doc/snort/README.wireless
share/doc/snort/generators
share/doc/snort/TODO
share/doc/snort/USAGE
share/doc/snort/WISHLIST
share/doc/snort/faq.pdf
share/doc/snort/faq.tex
share/doc/snort/schemas/create_db2
share/doc/snort/schemas/create_mssql
share/doc/snort/schemas/create_mysql
share/doc/snort/schemas/create_oracle.sql
share/doc/snort/schemas/create_postgresql
share/doc/snort/snort_manual.pdf
share/doc/snort/snort_manual.tex
share/doc/snort/snort_schema_v106.pdf
share/examples/rc.d/snort
share/examples/snort/classification.config
share/examples/snort/gen-msg.map
share/examples/snort/generators
share/examples/snort/reference.config
share/examples/snort/sid-msg.map
share/examples/snort/snort.conf.default
share/examples/snort/threshold.conf
share/examples/snort/unicode.map
share/snort/src/snort_dynamicsrc/bitop.h
share/snort/src/snort_dynamicsrc/debug.h
share/snort/src/snort_dynamicsrc/pcap_pkthdr32.h
share/snort/src/snort_dynamicsrc/preprocids.h
share/snort/src/snort_dynamicsrc/profiler.h
share/snort/src/snort_dynamicsrc/sf_dynamic_common.h
share/snort/src/snort_dynamicsrc/sf_dynamic_meta.h
share/snort/src/snort_dynamicsrc/sf_dynamic_preproc_lib.c
share/snort/src/snort_dynamicsrc/sf_dynamic_preproc_lib.h
share/snort/src/snort_dynamicsrc/sf_dynamic_preprocessor.h
share/snort/src/snort_dynamicsrc/sf_snort_packet.h
share/snort/src/snort_dynamicsrc/sf_snort_plugin_api.h
share/snort/src/snort_dynamicsrc/sfghash.h
share/snort/src/snort_dynamicsrc/sfhashfcn.h
share/snort/src/snort_dynamicsrc/sfsnort_dynamic_detection_lib.c
share/snort/src/snort_dynamicsrc/sfsnort_dynamic_detection_lib.h
share/snort/src/snort_dynamicsrc/str_search.h
share/snort/src/snort_dynamicsrc/stream_api.h
@dirrm share/snort/src/snort_dynamicsrc
@dirrm share/snort/src
@dirrm share/snort/rules
@dirrm share/snort
@dirrm share/examples/snort
@dirrm share/doc/snort/schemas
@dirrm share/doc/snort
@dirrm lib/snort_dynamicpreprocessor
@dirrm lib/snort_dynamicengine
@dirrm lib/snort_dynamicrules
