# net/nsd/PLIST
@comment $NetBSD: PLIST,v 1.11 2013/10/29 15:13:34 pettai Exp $
man/man5/nsd.conf.5
man/man8/nsd-checkconf.8
man/man8/nsd-control.8
man/man8/nsd.8
sbin/nsd
sbin/nsd-checkconf
sbin/nsd-control
sbin/nsd-control-setup
share/examples/nsd/nsd.conf
share/examples/rc.d/nsd
@pkgdir etc/nsd
