# net/nsd/PLIST
@comment $NetBSD: PLIST,v 1.8 2005/11/06 16:36:48 joerg Exp $
man/man8/nsd-notify.8
man/man8/nsd-xfer.8
man/man8/nsd.8
man/man8/nsdc.8
man/man8/zonec.8
sbin/nsd
sbin/nsd-notify
sbin/nsd-xfer
sbin/nsdc
sbin/zonec
share/examples/nsd/nsd.zones
share/examples/nsd/nsdc.conf
share/examples/rc.d/nsd
share/nsd/DIFFERENCES
share/nsd/LICENSE
share/nsd/NSD-DATABASE
share/nsd/README
share/nsd/RELNOTES
share/nsd/TODO
share/nsd/contrib/README
share/nsd/contrib/build-nsdzones.pl
share/nsd/contrib/build-nsdzones.pl.README
share/nsd/contrib/example-plugin.c
@dirrm share/nsd/contrib
@dirrm share/nsd
@dirrm share/examples/nsd
