# mail/sqwebmail/PLIST
@comment $NetBSD: PLIST,v 1.5 2005/05/02 20:34:01 reed Exp $
bin/mimegpg
libexec/cgi-bin/sqwebmail
libexec/courier/pcpd
libexec/courier/sqwebmaild
libexec/courier/sqwebpasswd
man/man1/mimegpg.1
sbin/sqwebmaild
sbin/webgpg
share/courier/sqwebmail/cleancache.pl
share/courier/sqwebmail/html/en
share/courier/sqwebmail/html/en-us/CHARSET
share/courier/sqwebmail/html/en-us/ISPELLDICT
share/courier/sqwebmail/html/en-us/LANGUAGE
share/courier/sqwebmail/html/en-us/LANGUAGE_PREF
share/courier/sqwebmail/html/en-us/LOCALE
share/courier/sqwebmail/html/en-us/TIMEZONELIST
share/courier/sqwebmail/html/en-us/abooklist.html
share/courier/sqwebmail/html/en-us/acl.html
share/courier/sqwebmail/html/en-us/attachments.html
share/courier/sqwebmail/html/en-us/autoresponder.html
share/courier/sqwebmail/html/en-us/calendarlogin.inc.html
share/courier/sqwebmail/html/en-us/empty.html
share/courier/sqwebmail/html/en-us/eventacl.html
share/courier/sqwebmail/html/en-us/eventdaily.html
share/courier/sqwebmail/html/en-us/eventdelete.html
share/courier/sqwebmail/html/en-us/eventmonthly.html
share/courier/sqwebmail/html/en-us/eventnotifydelete.txt
share/courier/sqwebmail/html/en-us/eventnotifynew.txt
share/courier/sqwebmail/html/en-us/eventnotifysubject.txt
share/courier/sqwebmail/html/en-us/eventshow.html
share/courier/sqwebmail/html/en-us/eventweekly.html
share/courier/sqwebmail/html/en-us/expired.html
share/courier/sqwebmail/html/en-us/filter.html
share/courier/sqwebmail/html/en-us/folder.html
share/courier/sqwebmail/html/en-us/folders.html
share/courier/sqwebmail/html/en-us/gpg.html
share/courier/sqwebmail/html/en-us/gpgcreate.html
share/courier/sqwebmail/html/en-us/gpgerr.html
share/courier/sqwebmail/html/en-us/index.html
share/courier/sqwebmail/html/en-us/invalid.html
share/courier/sqwebmail/html/en-us/keyimport.html
share/courier/sqwebmail/html/en-us/ldaplist.html
share/courier/sqwebmail/html/en-us/ldapsearch.html
share/courier/sqwebmail/html/en-us/login.html
share/courier/sqwebmail/html/en-us/navbar.inc.html
share/courier/sqwebmail/html/en-us/navbar2.inc.html
share/courier/sqwebmail/html/en-us/navbar3.inc.html
share/courier/sqwebmail/html/en-us/newevent.html
share/courier/sqwebmail/html/en-us/newmsg.html
share/courier/sqwebmail/html/en-us/preferences.html
share/courier/sqwebmail/html/en-us/print.html
share/courier/sqwebmail/html/en-us/printnocookie.html
share/courier/sqwebmail/html/en-us/printredirect.html
share/courier/sqwebmail/html/en-us/quickadd.html
share/courier/sqwebmail/html/en-us/readmsg.html
share/courier/sqwebmail/html/en-us/redirect.html
share/courier/sqwebmail/html/en-us/smileys.txt
share/courier/sqwebmail/html/en-us/spellchk.html
share/courier/sqwebmail/images/bg.png
share/courier/sqwebmail/images/bluebottomleft.png
share/courier/sqwebmail/images/bluebottomright.png
share/courier/sqwebmail/images/bluetopleft.png
share/courier/sqwebmail/images/bluetopright.png
share/courier/sqwebmail/images/calendar.png
share/courier/sqwebmail/images/cancel.png
share/courier/sqwebmail/images/eventadd.png
share/courier/sqwebmail/images/eventedit.png
share/courier/sqwebmail/images/first.png
share/courier/sqwebmail/images/folder.png
share/courier/sqwebmail/images/folder2.png
share/courier/sqwebmail/images/folders.png
share/courier/sqwebmail/images/forward.png
share/courier/sqwebmail/images/forwardatt.png
share/courier/sqwebmail/images/fullheaders.png
share/courier/sqwebmail/images/graybottomleft.png
share/courier/sqwebmail/images/graybottomright.png
share/courier/sqwebmail/images/graytopleft.png
share/courier/sqwebmail/images/graytopright.png
share/courier/sqwebmail/images/key.png
share/courier/sqwebmail/images/last.png
share/courier/sqwebmail/images/left.png
share/courier/sqwebmail/images/left2.png
share/courier/sqwebmail/images/logo.png
share/courier/sqwebmail/images/print.png
share/courier/sqwebmail/images/reply.png
share/courier/sqwebmail/images/replyall.png
share/courier/sqwebmail/images/replylist.png
share/courier/sqwebmail/images/right.png
share/courier/sqwebmail/images/right2.png
share/courier/sqwebmail/images/save.png
share/courier/sqwebmail/images/sep.png
share/courier/sqwebmail/images/signin.png
share/courier/sqwebmail/images/sm-frowney.png
share/courier/sqwebmail/images/sm-frowney2.png
share/courier/sqwebmail/images/sm-notsmiley.png
share/courier/sqwebmail/images/sm-notsmiley2.png
share/courier/sqwebmail/images/sm-ohno.png
share/courier/sqwebmail/images/sm-smiley.png
share/courier/sqwebmail/images/sm-smiley2.png
share/courier/sqwebmail/images/sm-wink.png
share/courier/sqwebmail/images/sm-wink2.png
share/courier/sqwebmail/images/sqwebmail.css
share/courier/sqwebmail/images/trash2.png
share/courier/sqwebmail/images/uncancel.png
share/courier/sqwebmail/ldapsearch
share/courier/sqwebmail/sendit.sh
share/doc/courier/sqwebmail/BUGS.html
share/doc/courier/sqwebmail/INSTALL.html
share/doc/courier/sqwebmail/README.gpg.html
share/doc/courier/sqwebmail/README.logindomainlist.html
share/doc/courier/sqwebmail/README.pcp.html
share/doc/courier/sqwebmail/SECURITY.html
share/examples/courier/calendarmode
share/examples/courier/ldapaddressbook.dist
share/examples/courier/sqwebmaild.dist
share/examples/courier/webmail.pam
share/examples/rc.d/sqwebmail
@unexec ${RMDIR} -p %D/share/examples/courier 2>/dev/null || ${TRUE}
@dirrm share/doc/courier/sqwebmail
@unexec ${RMDIR} -p %D/share/doc/html/courier 2>/dev/null || ${TRUE}
@dirrm share/courier/sqwebmail/images
@dirrm share/courier/sqwebmail/html/en-us
@dirrm share/courier/sqwebmail/html
@dirrm share/courier/sqwebmail
@unexec ${RMDIR} -p %D/share/courier 2>/dev/null || ${TRUE}
@unexec ${RMDIR} -p %D/libexec/courier 2>/dev/null || ${TRUE}
