# devel/p4d/PLIST
@comment $NetBSD: PLIST,v 1.3 2005/05/02 20:33:59 reed Exp $
sbin/p4d
share/examples/rc.d/p4d
