# devel/monotone-server/PLIST
@comment $NetBSD: PLIST,v 1.5 2009/06/14 17:48:59 joerg Exp $
sbin/monotone-server-init
share/examples/monotone-server/hooks.conf
share/examples/monotone-server/read-permissions
share/examples/monotone-server/write-permissions
share/examples/rc.d/monotone
