# filesystems/glusterfs8/PLIST
@comment $NetBSD: PLIST,v 1.1 2023/03/12 11:55:59 js Exp $
${PYSITELIB}/gluster/__init__.py
${PYSITELIB}/gluster/__init__.pyc
${PYSITELIB}/gluster/__init__.pyo
${PYSITELIB}/gluster/cliutils/__init__.py
${PYSITELIB}/gluster/cliutils/__init__.pyc
${PYSITELIB}/gluster/cliutils/__init__.pyo
${PYSITELIB}/gluster/cliutils/cliutils.py
${PYSITELIB}/gluster/cliutils/cliutils.pyc
${PYSITELIB}/gluster/cliutils/cliutils.pyo
bin/glusterfind
etc/ganesha/ganesha-ha.conf.sample
etc/glusterfs/eventsconfig.json
etc/glusterfs/gluster-rsyslog-5.8.conf
etc/glusterfs/gluster-rsyslog-7.2.conf
etc/glusterfs/glusterd.vol
etc/glusterfs/glusterfs-georep-logrotate
etc/glusterfs/glusterfs-logrotate
etc/glusterfs/group-db-workload
etc/glusterfs/group-distributed-virt
etc/glusterfs/group-gluster-block
etc/glusterfs/group-metadata-cache
etc/glusterfs/group-nl-cache
etc/glusterfs/group-samba
etc/glusterfs/group-virt.example
etc/glusterfs/gsyncd.conf
etc/glusterfs/logger.conf.example
etc/glusterfs/thin-arbiter.vol
include/glusterfs/api/glfs-handles.h
include/glusterfs/api/glfs.h
include/glusterfs/async.h
include/glusterfs/atomic.h
include/glusterfs/byte-order.h
include/glusterfs/call-stub.h
include/glusterfs/checksum.h
include/glusterfs/circ-buff.h
include/glusterfs/client_t.h
include/glusterfs/cluster-syncop.h
include/glusterfs/common-utils.h
include/glusterfs/compat-errno.h
include/glusterfs/compat-uuid.h
include/glusterfs/compat.h
include/glusterfs/daemon.h
include/glusterfs/default-args.h
include/glusterfs/defaults.h
include/glusterfs/dict.h
include/glusterfs/event-history.h
include/glusterfs/events.h
include/glusterfs/fd-lk.h
include/glusterfs/fd.h
include/glusterfs/gd-common-utils.h
include/glusterfs/gf-dirent.h
include/glusterfs/gf-event.h
include/glusterfs/gfchangelog/changelog.h
include/glusterfs/gidcache.h
include/glusterfs/glfs-message-id.h
include/glusterfs/globals.h
include/glusterfs/glusterfs-acl.h
include/glusterfs/glusterfs-fops.h
include/glusterfs/glusterfs.h
include/glusterfs/graph-utils.h
include/glusterfs/hashfn.h
include/glusterfs/iatt.h
include/glusterfs/inode.h
include/glusterfs/iobuf.h
include/glusterfs/latency.h
include/glusterfs/libglusterfs-messages.h
include/glusterfs/list.h
include/glusterfs/lkowner.h
include/glusterfs/locking.h
include/glusterfs/logging.h
include/glusterfs/lvm-defaults.h
include/glusterfs/mem-pool.h
include/glusterfs/mem-types.h
include/glusterfs/monitoring.h
include/glusterfs/options.h
include/glusterfs/parse-utils.h
include/glusterfs/quota-common-utils.h
include/glusterfs/rbthash.h
include/glusterfs/refcount.h
include/glusterfs/revision.h
include/glusterfs/rot-buffs.h
include/glusterfs/rpc/changelog-xdr.h
include/glusterfs/rpc/cli1-xdr.h
include/glusterfs/rpc/glusterd1-xdr.h
include/glusterfs/rpc/glusterfs3-xdr.h
include/glusterfs/rpc/glusterfs3.h
include/glusterfs/rpc/glusterfs4-xdr.h
include/glusterfs/rpc/portmap-xdr.h
include/glusterfs/rpc/protocol-common.h
include/glusterfs/rpc/rpc-clnt-ping.h
include/glusterfs/rpc/rpc-clnt.h
include/glusterfs/rpc/rpc-common-xdr.h
include/glusterfs/rpc/rpc-drc.h
include/glusterfs/rpc/rpc-lib-messages.h
include/glusterfs/rpc/rpc-pragmas.h
include/glusterfs/rpc/rpc-transport.h
include/glusterfs/rpc/rpcsvc-common.h
include/glusterfs/rpc/rpcsvc.h
include/glusterfs/rpc/xdr-common.h
include/glusterfs/rpc/xdr-generic.h
include/glusterfs/rpc/xdr-rpc.h
include/glusterfs/rpc/xdr-rpcclnt.h
include/glusterfs/run.h
include/glusterfs/server/authenticate.h
include/glusterfs/server/server-common.h
include/glusterfs/server/server-helpers.h
include/glusterfs/server/server-mem-types.h
include/glusterfs/server/server-messages.h
include/glusterfs/server/server.h
include/glusterfs/stack.h
include/glusterfs/statedump.h
include/glusterfs/store.h
include/glusterfs/strfd.h
include/glusterfs/syncop-utils.h
include/glusterfs/syncop.h
include/glusterfs/syscall.h
include/glusterfs/template-component-messages.h
include/glusterfs/throttle-tbf.h
include/glusterfs/timer.h
include/glusterfs/timespec.h
include/glusterfs/trie.h
include/glusterfs/upcall-utils.h
include/glusterfs/xlator.h
include/glusterfs/y.tab.h
lib/glusterfs/${GLUSTERFS_VERSION}/auth/addr.la
lib/glusterfs/${GLUSTERFS_VERSION}/auth/login.la
lib/glusterfs/${GLUSTERFS_VERSION}/rpc-transport/socket.la
lib/glusterfs/${GLUSTERFS_VERSION}/xlator/cluster/afr.la
lib/glusterfs/${GLUSTERFS_VERSION}/xlator/cluster/dht.la
lib/glusterfs/${GLUSTERFS_VERSION}/xlator/cluster/disperse.so
lib/glusterfs/${GLUSTERFS_VERSION}/xlator/cluster/distribute.so
lib/glusterfs/${GLUSTERFS_VERSION}/xlator/cluster/ec.la
lib/glusterfs/${GLUSTERFS_VERSION}/xlator/cluster/nufa.la
lib/glusterfs/${GLUSTERFS_VERSION}/xlator/cluster/replicate.so
lib/glusterfs/${GLUSTERFS_VERSION}/xlator/cluster/switch.la
lib/glusterfs/${GLUSTERFS_VERSION}/xlator/debug/delay-gen.la
lib/glusterfs/${GLUSTERFS_VERSION}/xlator/debug/error-gen.la
lib/glusterfs/${GLUSTERFS_VERSION}/xlator/debug/io-stats.la
lib/glusterfs/${GLUSTERFS_VERSION}/xlator/debug/sink.la
lib/glusterfs/${GLUSTERFS_VERSION}/xlator/debug/trace.la
lib/glusterfs/${GLUSTERFS_VERSION}/xlator/features/access-control.so
lib/glusterfs/${GLUSTERFS_VERSION}/xlator/features/arbiter.la
lib/glusterfs/${GLUSTERFS_VERSION}/xlator/features/barrier.la
lib/glusterfs/${GLUSTERFS_VERSION}/xlator/features/bit-rot.la
lib/glusterfs/${GLUSTERFS_VERSION}/xlator/features/bitrot-stub.la
lib/glusterfs/${GLUSTERFS_VERSION}/xlator/features/cdc.la
lib/glusterfs/${GLUSTERFS_VERSION}/xlator/features/changelog.la
lib/glusterfs/${GLUSTERFS_VERSION}/xlator/features/gfid-access.la
lib/glusterfs/${GLUSTERFS_VERSION}/xlator/features/index.la
lib/glusterfs/${GLUSTERFS_VERSION}/xlator/features/leases.la
lib/glusterfs/${GLUSTERFS_VERSION}/xlator/features/locks.la
lib/glusterfs/${GLUSTERFS_VERSION}/xlator/features/marker.la
lib/glusterfs/${GLUSTERFS_VERSION}/xlator/features/namespace.la
lib/glusterfs/${GLUSTERFS_VERSION}/xlator/features/posix-locks.so
lib/glusterfs/${GLUSTERFS_VERSION}/xlator/features/quiesce.la
lib/glusterfs/${GLUSTERFS_VERSION}/xlator/features/quota.la
lib/glusterfs/${GLUSTERFS_VERSION}/xlator/features/quotad.la
lib/glusterfs/${GLUSTERFS_VERSION}/xlator/features/read-only.la
lib/glusterfs/${GLUSTERFS_VERSION}/xlator/features/sdfs.la
lib/glusterfs/${GLUSTERFS_VERSION}/xlator/features/selinux.la
lib/glusterfs/${GLUSTERFS_VERSION}/xlator/features/shard.la
lib/glusterfs/${GLUSTERFS_VERSION}/xlator/features/snapview-client.la
lib/glusterfs/${GLUSTERFS_VERSION}/xlator/features/snapview-server.la
lib/glusterfs/${GLUSTERFS_VERSION}/xlator/features/thin-arbiter.la
lib/glusterfs/${GLUSTERFS_VERSION}/xlator/features/trash.la
lib/glusterfs/${GLUSTERFS_VERSION}/xlator/features/upcall.la
lib/glusterfs/${GLUSTERFS_VERSION}/xlator/features/utime.la
lib/glusterfs/${GLUSTERFS_VERSION}/xlator/features/worm.la
lib/glusterfs/${GLUSTERFS_VERSION}/xlator/meta.la
lib/glusterfs/${GLUSTERFS_VERSION}/xlator/mgmt/glusterd.la
lib/glusterfs/${GLUSTERFS_VERSION}/xlator/mount/api.la
lib/glusterfs/${GLUSTERFS_VERSION}/xlator/mount/fuse.la
lib/glusterfs/${GLUSTERFS_VERSION}/xlator/performance/io-cache.la
lib/glusterfs/${GLUSTERFS_VERSION}/xlator/performance/io-threads.la
lib/glusterfs/${GLUSTERFS_VERSION}/xlator/performance/md-cache.la
lib/glusterfs/${GLUSTERFS_VERSION}/xlator/performance/nl-cache.la
lib/glusterfs/${GLUSTERFS_VERSION}/xlator/performance/open-behind.la
lib/glusterfs/${GLUSTERFS_VERSION}/xlator/performance/quick-read.la
lib/glusterfs/${GLUSTERFS_VERSION}/xlator/performance/read-ahead.la
lib/glusterfs/${GLUSTERFS_VERSION}/xlator/performance/readdir-ahead.la
lib/glusterfs/${GLUSTERFS_VERSION}/xlator/performance/stat-prefetch.so
lib/glusterfs/${GLUSTERFS_VERSION}/xlator/performance/write-behind.la
lib/glusterfs/${GLUSTERFS_VERSION}/xlator/playground/template.la
lib/glusterfs/${GLUSTERFS_VERSION}/xlator/protocol/client.la
lib/glusterfs/${GLUSTERFS_VERSION}/xlator/protocol/server.la
lib/glusterfs/${GLUSTERFS_VERSION}/xlator/storage/posix.la
lib/glusterfs/${GLUSTERFS_VERSION}/xlator/system/posix-acl.la
lib/libgfapi.la
lib/libgfchangelog.la
lib/libgfrpc.la
lib/libgfxdr.la
lib/libglusterd.la
lib/libglusterfs.la
lib/ocf/resource.d/glusterfs/glusterd
lib/ocf/resource.d/glusterfs/volume
lib/ocf/resource.d/heartbeat/ganesha_grace
lib/ocf/resource.d/heartbeat/ganesha_mon
lib/ocf/resource.d/heartbeat/ganesha_nfsd
lib/pkgconfig/glusterfs-api.pc
lib/pkgconfig/libgfchangelog.pc
libexec/ganesha/create-export-ganesha.sh
libexec/ganesha/dbus-send.sh
libexec/ganesha/ganesha-ha.sh
libexec/ganesha/generate-epoch.py
libexec/glusterfs/gfevents/__init__.py
libexec/glusterfs/gfevents/__init__.pyc
libexec/glusterfs/gfevents/__init__.pyo
libexec/glusterfs/gfevents/eventsapiconf.py
libexec/glusterfs/gfevents/eventsapiconf.pyc
libexec/glusterfs/gfevents/eventsapiconf.pyo
libexec/glusterfs/gfevents/eventtypes.py
libexec/glusterfs/gfevents/eventtypes.pyc
libexec/glusterfs/gfevents/eventtypes.pyo
libexec/glusterfs/gfevents/gf_event.py
libexec/glusterfs/gfevents/gf_event.pyc
libexec/glusterfs/gfevents/gf_event.pyo
libexec/glusterfs/gfevents/glustereventsd.py
libexec/glusterfs/gfevents/handlers.py
libexec/glusterfs/gfevents/handlers.pyc
libexec/glusterfs/gfevents/handlers.pyo
libexec/glusterfs/gfevents/utils.py
libexec/glusterfs/gfevents/utils.pyc
libexec/glusterfs/gfevents/utils.pyo
libexec/glusterfs/gfind_missing_files/gcrawler
libexec/glusterfs/gfind_missing_files/gfid_to_path.py
libexec/glusterfs/gfind_missing_files/gfid_to_path.sh
libexec/glusterfs/gfind_missing_files/gfind_missing_files.sh
libexec/glusterfs/glfsheal
libexec/glusterfs/glusterfind/S57glusterfind-delete-post.py
libexec/glusterfs/glusterfind/__init__.py
libexec/glusterfs/glusterfind/__init__.pyc
libexec/glusterfs/glusterfind/__init__.pyo
libexec/glusterfs/glusterfind/brickfind.py
libexec/glusterfs/glusterfind/changelog.py
libexec/glusterfs/glusterfind/changelogdata.py
libexec/glusterfs/glusterfind/changelogdata.pyc
libexec/glusterfs/glusterfind/changelogdata.pyo
libexec/glusterfs/glusterfind/conf.py
libexec/glusterfs/glusterfind/conf.pyc
libexec/glusterfs/glusterfind/conf.pyo
libexec/glusterfs/glusterfind/gfind_py2py3.py
libexec/glusterfs/glusterfind/gfind_py2py3.pyc
libexec/glusterfs/glusterfind/gfind_py2py3.pyo
libexec/glusterfs/glusterfind/libgfchangelog.py
libexec/glusterfs/glusterfind/libgfchangelog.pyc
libexec/glusterfs/glusterfind/libgfchangelog.pyo
libexec/glusterfs/glusterfind/main.py
libexec/glusterfs/glusterfind/main.pyc
libexec/glusterfs/glusterfind/main.pyo
libexec/glusterfs/glusterfind/nodeagent.py
libexec/glusterfs/glusterfind/tool.conf
libexec/glusterfs/glusterfind/utils.py
libexec/glusterfs/glusterfind/utils.pyc
libexec/glusterfs/glusterfind/utils.pyo
libexec/glusterfs/gsyncd
libexec/glusterfs/gverify.sh
libexec/glusterfs/peer_add_secret_pub
libexec/glusterfs/peer_eventsapi.py
libexec/glusterfs/peer_georep-sshkey.py
libexec/glusterfs/peer_gsec_create
libexec/glusterfs/peer_mountbroker
libexec/glusterfs/peer_mountbroker.py
libexec/glusterfs/python/syncdaemon/README.md
libexec/glusterfs/python/syncdaemon/__init__.py
libexec/glusterfs/python/syncdaemon/__init__.pyc
libexec/glusterfs/python/syncdaemon/__init__.pyo
libexec/glusterfs/python/syncdaemon/argsupgrade.py
libexec/glusterfs/python/syncdaemon/argsupgrade.pyc
libexec/glusterfs/python/syncdaemon/argsupgrade.pyo
libexec/glusterfs/python/syncdaemon/conf.py
libexec/glusterfs/python/syncdaemon/conf.pyc
libexec/glusterfs/python/syncdaemon/conf.pyo
libexec/glusterfs/python/syncdaemon/gsyncd.py
libexec/glusterfs/python/syncdaemon/gsyncd.pyc
libexec/glusterfs/python/syncdaemon/gsyncd.pyo
libexec/glusterfs/python/syncdaemon/gsyncdconfig.py
libexec/glusterfs/python/syncdaemon/gsyncdconfig.pyc
libexec/glusterfs/python/syncdaemon/gsyncdconfig.pyo
libexec/glusterfs/python/syncdaemon/gsyncdstatus.py
libexec/glusterfs/python/syncdaemon/gsyncdstatus.pyc
libexec/glusterfs/python/syncdaemon/gsyncdstatus.pyo
libexec/glusterfs/python/syncdaemon/libcxattr.py
libexec/glusterfs/python/syncdaemon/libcxattr.pyc
libexec/glusterfs/python/syncdaemon/libcxattr.pyo
libexec/glusterfs/python/syncdaemon/libgfchangelog.py
libexec/glusterfs/python/syncdaemon/libgfchangelog.pyc
libexec/glusterfs/python/syncdaemon/libgfchangelog.pyo
libexec/glusterfs/python/syncdaemon/logutils.py
libexec/glusterfs/python/syncdaemon/logutils.pyc
libexec/glusterfs/python/syncdaemon/logutils.pyo
libexec/glusterfs/python/syncdaemon/master.py
libexec/glusterfs/python/syncdaemon/master.pyc
libexec/glusterfs/python/syncdaemon/master.pyo
libexec/glusterfs/python/syncdaemon/monitor.py
libexec/glusterfs/python/syncdaemon/monitor.pyc
libexec/glusterfs/python/syncdaemon/monitor.pyo
libexec/glusterfs/python/syncdaemon/py2py3.py
libexec/glusterfs/python/syncdaemon/py2py3.pyc
libexec/glusterfs/python/syncdaemon/py2py3.pyo
libexec/glusterfs/python/syncdaemon/rconf.py
libexec/glusterfs/python/syncdaemon/rconf.pyc
libexec/glusterfs/python/syncdaemon/rconf.pyo
libexec/glusterfs/python/syncdaemon/repce.py
libexec/glusterfs/python/syncdaemon/repce.pyc
libexec/glusterfs/python/syncdaemon/repce.pyo
libexec/glusterfs/python/syncdaemon/resource.py
libexec/glusterfs/python/syncdaemon/resource.pyc
libexec/glusterfs/python/syncdaemon/resource.pyo
libexec/glusterfs/python/syncdaemon/subcmds.py
libexec/glusterfs/python/syncdaemon/subcmds.pyc
libexec/glusterfs/python/syncdaemon/subcmds.pyo
libexec/glusterfs/python/syncdaemon/syncdutils.py
libexec/glusterfs/python/syncdaemon/syncdutils.pyc
libexec/glusterfs/python/syncdaemon/syncdutils.pyo
libexec/glusterfs/scripts/generate-gfid-file.sh
libexec/glusterfs/scripts/get-gfid.sh
libexec/glusterfs/scripts/gsync-sync-gfid
libexec/glusterfs/scripts/gsync-upgrade.sh
libexec/glusterfs/scripts/schedule_georep.py
libexec/glusterfs/scripts/slave-upgrade.sh
libexec/glusterfs/set_geo_rep_pem_keys.sh
man/man8/gluster-setgfid2path.8
man/man8/gluster.8
man/man8/glusterd.8
man/man8/glusterfs.8
man/man8/glusterfsd.8
man/man8/mount.glusterfs.8
sbin/conf.py
sbin/gcron.py
sbin/gf_attach
sbin/gfind_missing_files
sbin/gluster
sbin/gluster-eventsapi
sbin/gluster-georep-sshkey
sbin/gluster-mountbroker
sbin/gluster-setgfid2path
sbin/glusterd
sbin/glustereventsd
sbin/glusterfs
sbin/glusterfsd
sbin/mount_glusterfs
sbin/snap_scheduler.py
sbin/umountd
share/doc/glusterfs/benchmarking/README
share/doc/glusterfs/benchmarking/glfs-bm.c
share/doc/glusterfs/benchmarking/launch-script.sh
share/doc/glusterfs/benchmarking/local-script.sh
share/doc/glusterfs/benchmarking/rdd.c
share/doc/glusterfs/glusterfs-mode.el
share/doc/glusterfs/glusterfs.vim
share/examples/rc.d/glusterd
share/glusterfs/scripts/eventsdash.py
share/glusterfs/scripts/post-upgrade-script-for-quota.sh
share/glusterfs/scripts/pre-upgrade-script-for-quota.sh
share/glusterfs/scripts/setup-thin-arbiter.sh
share/glusterfs/scripts/stop-all-gluster-processes.sh
