# mail/dspam/PLIST
@comment $NetBSD: PLIST,v 1.17 2007/05/07 09:29:54 adrianp Exp $
${HASH}bin/cssstat
${HASH}bin/csscompress
${HASH}bin/cssclean
${HASH}bin/cssconvert
bin/dspam
bin/dspam_2sql
bin/dspam_admin
bin/dspam_clean
bin/dspam_crc
bin/dspam_dump
bin/dspam_logrotate
bin/dspam_merge
bin/dspam_stats
bin/dspam_train
bin/dspamc
${PGSQL}bin/dspam_pg2int8
include/dspam/buffer.h
include/dspam/config.h
include/dspam/config_shared.h
include/dspam/decode.h
include/dspam/diction.h
include/dspam/error.h
include/dspam/heap.h
include/dspam/ldap_client.h
include/dspam/libdspam.h
include/dspam/libdspam_objects.h
include/dspam/nodetree.h
include/dspam/storage_driver.h
include/dspam/tokenizer.h
lib/libdspam.la
lib/pkgconfig/dspam.pc
man/man1/dspam.1
man/man1/dspam_clean.1
man/man1/dspam_dump.1
man/man1/dspam_merge.1
man/man1/dspam_stats.1
man/man1/dspam_train.1
man/man3/dspam_addattribute.3
man/man3/dspam_attach.3
man/man3/dspam_create.3
man/man3/dspam_destroy.3
man/man3/dspam_detach.3
man/man3/dspam_getsource.3
man/man3/dspam_init.3
man/man3/dspam_process.3
man/man3/libdspam.3
share/doc/dspam/README
share/doc/dspam/RELEASE.NOTES
share/doc/dspam/UPGRADING
share/doc/dspam/courier.txt
share/doc/dspam/exim.txt
share/doc/dspam/markov.txt
share/doc/dspam/mysql_drv.txt
share/doc/dspam/osx.txt
share/doc/dspam/pgsql_drv.txt
share/doc/dspam/pop3filter.txt
share/doc/dspam/postfix.txt
share/doc/dspam/qmail.txt
share/doc/dspam/relay.txt
share/doc/dspam/sendmail.txt
share/doc/dspam/sqlite_drv.txt
share/doc/dspam/tests.txt
share/dspam/txt/firstrun.txt
share/dspam/txt/firstspam.txt
share/dspam/txt/quarantinefull.txt
share/dspam/webui/cgi-bin/admin.cgi
share/dspam/webui/cgi-bin/admingraph.cgi
share/dspam/webui/cgi-bin/dspam.cgi
share/dspam/webui/cgi-bin/graph.cgi
share/dspam/webui/cgi-bin/rgb.txt
share/dspam/webui/cgi-bin/templates/nav_admin_error.html
share/dspam/webui/cgi-bin/templates/nav_admin_preferences.html
share/dspam/webui/cgi-bin/templates/nav_admin_status.html
share/dspam/webui/cgi-bin/templates/nav_admin_user.html
share/dspam/webui/cgi-bin/templates/nav_alerts.html
share/dspam/webui/cgi-bin/templates/nav_analysis.html
share/dspam/webui/cgi-bin/templates/nav_error.html
share/dspam/webui/cgi-bin/templates/nav_fragment.html
share/dspam/webui/cgi-bin/templates/nav_history.html
share/dspam/webui/cgi-bin/templates/nav_performance.html
share/dspam/webui/cgi-bin/templates/nav_preferences.html
share/dspam/webui/cgi-bin/templates/nav_quarantine.html
share/dspam/webui/cgi-bin/templates/nav_viewmessage.html
share/dspam/webui/htdocs/base.css
share/dspam/webui/htdocs/dspam-logo-small.gif
${MYSQL}share/dspam/mysql/mysql_objects-4.1.sql
${MYSQL}share/dspam/mysql/mysql_objects-space.sql
${MYSQL}share/dspam/mysql/mysql_objects-speed.sql
${MYSQL}share/dspam/mysql/purge.sql
${MYSQL}share/dspam/mysql/purge-4.1.sql
${MYSQL}share/dspam/mysql/virtual_users.sql
${MYSQL}share/dspam/mysql/virtual_user_aliases.sql
${PGSQL}share/dspam/pgsql/pgsql_objects.sql
${PGSQL}share/dspam/pgsql/purge.sql
${PGSQL}share/dspam/pgsql/virtual_users.sql
${SQLITE}share/dspam/sqlite/purge-2.sql
${SQLITE}share/dspam/sqlite/purge-3.sql
share/examples/dspam/cgi-admins
share/examples/dspam/cgi-default.prefs
share/examples/dspam/configure.pl
share/examples/dspam/dspam.conf
share/examples/rc.d/dspam
@dirrm share/examples/dspam
${PGSQL}@dirrm share/dspam/pgsql
${MYSQL}@dirrm share/dspam/mysql
${SQLITE}@dirrm share/dspam/sqlite
@dirrm share/dspam/webui/htdocs
@dirrm share/dspam/webui/cgi-bin/templates
@dirrm share/dspam/webui/cgi-bin
@dirrm share/dspam/webui
@dirrm share/dspam/txt
@dirrm share/dspam
@dirrm share/doc/dspam
@dirrm include/dspam
