# sysutils/rsyslog/PLIST
@comment $NetBSD: PLIST,v 1.2 2014/12/12 12:47:46 fhajny Exp $
lib/rsyslog/imdiag.la
lib/rsyslog/imfile.la
${PLIST.klog}lib/rsyslog/imklog.la
lib/rsyslog/immark.la
lib/rsyslog/impstats.la
${PLIST.solaris}lib/rsyslog/imsolaris.la
lib/rsyslog/imtcp.la
lib/rsyslog/imudp.la
lib/rsyslog/imuxsock.la
lib/rsyslog/lmnet.la
lib/rsyslog/lmnetstrms.la
lib/rsyslog/lmnsd_ptcp.la
lib/rsyslog/lmregexp.la
lib/rsyslog/lmstrmsrv.la
lib/rsyslog/lmtcpclt.la
lib/rsyslog/lmtcpsrv.la
lib/rsyslog/lmzlibw.la
lib/rsyslog/mmexternal.la
lib/rsyslog/mmjsonparse.la
lib/rsyslog/mmnormalize.la
lib/rsyslog/ommail.la
lib/rsyslog/omstdout.la
lib/rsyslog/omtesting.la
man/man5/rsyslog.conf.5
man/man8/rsyslogd.8
sbin/msggen
sbin/rsyslog_diag_hostname
sbin/rsyslogd
share/examples/rc.d/rsyslogd
share/examples/rsyslog/rsyslog.conf
