# net/tor/PLIST
@comment $NetBSD: PLIST,v 1.5 2005/08/09 15:16:39 salo Exp $
bin/tor
bin/tor-resolve
bin/torify
man/man1/tor-resolve.1
man/man1/tor.1
man/man1/torify.1
share/examples/rc.d/tor
share/examples/tor/tor-tsocks.conf
share/examples/tor/torrc.sample
@dirrm share/examples/tor
