# security/nessus-core/PLIST
@comment $NetBSD: PLIST,v 1.15 2005/05/02 20:34:05 reed Exp $
bin/nessus
bin/nessus-fetch
bin/nessus-mkcert-client
bin/nessus-mkrand
include/nessus/config.h
include/nessus/includes.h
include/nessus/nessus-devel.h
include/nessus/nessusicmp.h
include/nessus/nessusip.h
include/nessus/nessusraw.h
include/nessus/nessustcp.h
include/nessus/nessusudp.h
include/nessus/ntcompat.h
man/man1/nessus-fetch.1
man/man1/nessus-mkcert-client.1
man/man1/nessus-mkrand.1
man/man1/nessus.1
man/man8/nessus-adduser.8
man/man8/nessus-mkcert.8
man/man8/nessus-rmuser.8
man/man8/nessusd.8
sbin/nessus-adduser
sbin/nessus-check-signature
sbin/nessus-mkcert
sbin/nessus-rmuser
sbin/nessusd
share/examples/rc.d/nessusd
@exec ${MKDIR} %D/lib/nessus/reports
@unexec ${RMDIR} %D/lib/nessus/reports 2>/dev/null || ${TRUE}
@unexec ${RMDIR} %D/lib/nessus 2>/dev/null || ${TRUE}
@exec ${MKDIR} %D/etc/nessus
@unexec ${RMDIR} %D/etc/nessus 2>/dev/null || ${TRUE}
@exec ${MKDIR} ${VARBASE}/nessus
@unexec ${RMDIR} ${VARBASE}/nessus 2>/dev/null || ${TRUE}
