# security/openssh/PLIST
@comment $NetBSD: PLIST,v 1.3 2002/06/25 02:13:12 itojun Exp $
bin/scp
bin/sftp
bin/slogin
bin/ssh
bin/ssh-add
bin/ssh-agent
bin/ssh-keygen
bin/ssh-keyscan
etc/rc.d/sshd
libexec/sftp-server
libexec/ssh-keysign
man/man1/scp.1
man/man1/sftp.1
man/man1/slogin.1
man/man1/ssh-add.1
man/man1/ssh-agent.1
man/man1/ssh-keygen.1
man/man1/ssh-keyscan.1
man/man1/ssh.1
man/man5/ssh_config.5
man/man5/sshd_config.5
man/man8/sftp-server.8
man/man8/ssh-keysign.8
man/man8/sshd.8
sbin/sshd
share/examples/openssh/moduli
share/examples/openssh/ssh_config
share/examples/openssh/sshd_config
@dirrm share/examples/openssh
@exec mkdir -p /var/chroot/sshd
