# news/inn/PLIST
@comment $NetBSD: PLIST,v 1.8 2005/05/02 20:34:04 reed Exp $
bin/inews
include/inn/config.h
include/inn/dbz.h
include/inn/libinn.h
include/inn/storage.h
inn/bin/actmerge
inn/bin/actsync
inn/bin/actsyncd
inn/bin/archive
inn/bin/auth/passwd/auth_smb
inn/bin/auth/passwd/ckpasswd
inn/bin/auth/passwd/radius
inn/bin/auth/resolv/domain
inn/bin/auth/resolv/ident
inn/bin/batcher
inn/bin/buffchan
inn/bin/cnfsheadconf
inn/bin/cnfsstat
inn/bin/control/checkgroups.pl
inn/bin/control/ihave.pl
inn/bin/control/newgroup.pl
inn/bin/control/rmgroup.pl
inn/bin/control/sendme.pl
inn/bin/control/sendsys.pl
inn/bin/control/senduuname.pl
inn/bin/control/version.pl
inn/bin/controlbatch
inn/bin/controlchan
inn/bin/convdate
inn/bin/ctlinnd
inn/bin/cvtbatch
inn/bin/docheckgroups
inn/bin/expire
inn/bin/expireover
inn/bin/expirerm
inn/bin/fastrm
inn/bin/filechan
inn/bin/filter/INN.py
inn/bin/filter/filter.tcl
inn/bin/filter/filter_innd.pl
inn/bin/filter/filter_innd.py
inn/bin/filter/filter_nnrpd.pl
inn/bin/filter/nnrpd_auth.pl
inn/bin/filter/nnrpd_auth.py
inn/bin/filter/startup.tcl
inn/bin/filter/startup_innd.pl
inn/bin/getlist
inn/bin/gpgverify
inn/bin/grephistory
inn/bin/imapfeed
inn/bin/inews
inn/bin/inncheck
inn/bin/innconfval
inn/bin/innd
inn/bin/inndf
inn/bin/inndstart
inn/bin/innfeed
inn/bin/innmail
inn/bin/innreport
inn/bin/innstat
inn/bin/innupgrade
inn/bin/innwatch
inn/bin/innxbatch
inn/bin/innxmit
inn/bin/mailpost
inn/bin/makedbz
inn/bin/makehistory
inn/bin/mod-active
inn/bin/news.daily
inn/bin/news2mail
inn/bin/ninpaths
inn/bin/nnrpd
inn/bin/nntpget
inn/bin/nntpsend
inn/bin/ovdb_init
inn/bin/ovdb_monitor
inn/bin/ovdb_server
inn/bin/ovdb_stat
inn/bin/overchan
inn/bin/perl-nocem
inn/bin/pgpverify
inn/bin/procbatch
inn/bin/prunehistory
inn/bin/pullnews
inn/bin/rc.news
inn/bin/rnews
inn/bin/rnews.libexec/c7unbatch
inn/bin/rnews.libexec/decode
inn/bin/rnews.libexec/encode
inn/bin/rnews.libexec/gunbatch
inn/bin/scanlogs
inn/bin/scanspool
inn/bin/send-ihave
inn/bin/send-nntp
inn/bin/send-uucp
inn/bin/sendinpaths
inn/bin/sendxbatches
inn/bin/shlock
inn/bin/shrinkfile
inn/bin/signcontrol
inn/bin/simpleftp
inn/bin/sm
inn/bin/startinnfeed
inn/bin/tally.control
inn/bin/tdx-util
inn/bin/writelog
inn/doc/CONTRIBUTORS
inn/doc/GPL
inn/doc/HACKING
inn/doc/INSTALL
inn/doc/LICENSE
inn/doc/NEWS
inn/doc/README
inn/doc/TODO
inn/doc/compliance-nntp
inn/doc/config-design
inn/doc/config-semantics
inn/doc/config-syntax
inn/doc/external-auth
inn/doc/history
inn/doc/hook-perl
inn/doc/hook-python
inn/doc/hook-tcl
inn/doc/sample-control
inn/include/conffile.h
inn/include/config.h
inn/include/dbz.h
inn/include/inn/buffer.h
inn/include/inn/confparse.h
inn/include/inn/defines.h
inn/include/inn/hashtab.h
inn/include/inn/history.h
inn/include/inn/innconf.h
inn/include/inn/list.h
inn/include/inn/md5.h
inn/include/inn/messages.h
inn/include/inn/mmap.h
inn/include/inn/qio.h
inn/include/inn/sequence.h
inn/include/inn/system.h
inn/include/inn/timer.h
inn/include/inn/tst.h
inn/include/inn/vector.h
inn/include/inn/version.h
inn/include/inn/wire.h
inn/include/inndcomm.h
inn/include/libinn.h
inn/include/nntp.h
inn/include/ov.h
inn/include/paths.h
inn/include/storage.h
inn/lib/innreport_inn.pm
inn/lib/innshellvars
inn/lib/innshellvars.pl
inn/lib/innshellvars.tcl
inn/lib/libinn.a
inn/lib/libinnhist.a
inn/lib/libstorage.a
man/man1/convdate.1
man/man1/fastrm.1
man/man1/getlist.1
man/man1/grephistory.1
man/man1/inews.1
man/man1/innconfval.1
man/man1/innfeed.1
man/man1/innmail.1
man/man1/nntpget.1
man/man1/rnews.1
man/man1/shlock.1
man/man1/shrinkfile.1
man/man1/simpleftp.1
man/man1/startinnfeed.1
man/man3/clientlib.3
man/man3/dbz.3
man/man3/inndcomm.3
man/man3/libinn.3
man/man3/libstorage.3
man/man3/parsedate.3
man/man3/qio.3
man/man5/active.5
man/man5/control.ctl.5
man/man5/cycbuff.conf.5
man/man5/distrib.pats.5
man/man5/expire.ctl.5
man/man5/history.5
man/man5/incoming.conf.5
man/man5/inn.conf.5
man/man5/innfeed.conf.5
man/man5/innwatch.ctl.5
man/man5/moderators.5
man/man5/motd.news.5
man/man5/newsfeeds.5
man/man5/nnrpd.track.5
man/man5/nntpsend.ctl.5
man/man5/overview.fmt.5
man/man5/passwd.nntp.5
man/man5/storage.conf.5
man/man5/subscriptions.5
man/man8/actsync.8
man/man8/actsyncd.8
man/man8/archive.8
man/man8/batcher.8
man/man8/buffchan.8
man/man8/cnfsstat.8
man/man8/controlchan.8
man/man8/ctlinnd.8
man/man8/cvtbatch.8
man/man8/expire.8
man/man8/expireover.8
man/man8/expirerm.8
man/man8/filechan.8
man/man8/inncheck.8
man/man8/innd.8
man/man8/inndf.8
man/man8/innreport.8
man/man8/innstat.8
man/man8/innwatch.8
man/man8/innxbatch.8
man/man8/innxmit.8
man/man8/mailpost.8
man/man8/makehistory.8
man/man8/mod-active.8
man/man8/news.daily.8
man/man8/news2mail.8
man/man8/nnrpd.8
man/man8/nntpsend.8
man/man8/overchan.8
man/man8/prunehistory.8
man/man8/pullnews.8
man/man8/scanlogs.8
man/man8/send-uucp.8
man/man8/tally.control.8
man/man8/writelog.8
share/examples/inn/INN.py
share/examples/inn/active.minimal
share/examples/inn/actsync.cfg
share/examples/inn/actsync.ign
share/examples/inn/buffindexed.conf
share/examples/inn/control.ctl
share/examples/inn/cycbuff.conf
share/examples/inn/distrib.pats
share/examples/inn/expire.ctl
share/examples/inn/filter.tcl
share/examples/inn/filter_innd.pl
share/examples/inn/filter_innd.py
share/examples/inn/filter_nnrpd.pl
share/examples/inn/incoming.conf
share/examples/inn/inn.conf
share/examples/inn/innfeed.conf
share/examples/inn/innreport.conf
share/examples/inn/innwatch.ctl
share/examples/inn/moderators
share/examples/inn/motd.news
share/examples/inn/news2mail.cf
share/examples/inn/newsfeeds
share/examples/inn/newsgroups.minimal
share/examples/inn/nnrpd.py
share/examples/inn/nnrpd.track
share/examples/inn/nnrpd_access.pl
share/examples/inn/nnrpd_access_wrapper.pl
share/examples/inn/nnrpd_access_wrapper.py
share/examples/inn/nnrpd_auth.pl
share/examples/inn/nnrpd_auth.py
share/examples/inn/nnrpd_auth_wrapper.pl
share/examples/inn/nnrpd_auth_wrapper.py
share/examples/inn/nnrpd_dynamic_wrapper.py
share/examples/inn/nntpsend.ctl
share/examples/inn/ovdb.conf
share/examples/inn/overview.fmt
share/examples/inn/passwd.nntp
share/examples/inn/radius.conf
share/examples/inn/readers.conf
share/examples/inn/sasl.conf
share/examples/inn/startup.tcl
share/examples/inn/startup_innd.pl
share/examples/inn/storage.conf
share/examples/inn/subscriptions
share/examples/rc.d/innd
@dirrm share/examples/inn
@dirrm inn/lib
@dirrm inn/include/inn
@dirrm inn/include
@dirrm inn/doc
@dirrm inn/bin/rnews.libexec
@dirrm inn/bin/filter
@dirrm inn/bin/control
@dirrm inn/bin/auth/resolv
@dirrm inn/bin/auth/passwd
@dirrm inn/bin/auth
@dirrm inn/bin
@dirrm inn
@dirrm include/inn
