# sysutils/monit/PLIST
@comment $NetBSD: PLIST,v 1.2 2007/12/15 16:41:50 seb Exp $
bin/monit
man/man1/monit.1
share/docs/monit/CHANGES.txt
share/docs/monit/CONTRIBUTORS
share/docs/monit/COPYING
share/docs/monit/FAQ.txt
share/docs/monit/LICENSE
share/docs/monit/PACKAGES
share/docs/monit/PLATFORMS
share/docs/monit/README
share/docs/monit/README.SSL
share/docs/monit/STATUS
share/docs/monit/UPGRADE.txt
share/examples/monit/MacOSX/README
share/examples/monit/MacOSX/StartupParameters.plist
share/examples/monit/MacOSX/monit
share/examples/monit/MacOSX/monit.plist
share/examples/monit/MacOSX/monitrc
share/examples/monit/monit.php
share/examples/monit/monitrc-gentoo
share/examples/monit/monitrc.sample
share/examples/monit/wap.php
share/examples/rc.d/monit
@dirrm share/examples/monit/MacOSX
@dirrm share/examples/monit
@dirrm share/docs/monit
