# sysutils/monit/PLIST
@comment $NetBSD: PLIST,v 1.4 2014/02/06 09:09:17 obache Exp $
bin/monit
man/man1/monit.1
share/docs/monit/CHANGES
share/docs/monit/COPYING
share/docs/monit/PLATFORMS
share/docs/monit/README
share/examples/monit/monit.service
share/examples/monit/monit.upstart
share/examples/monit/monitrc-gentoo
share/examples/monit/monitrc.sample
share/examples/monit/rc.monit
share/examples/rc.d/monit
