# net/xymon/PLIST
@comment $NetBSD: PLIST,v 1.1.1.1 2010/02/14 22:58:57 spz Exp $
libexec/xymon/bb
libexec/xymon/bb-ack.cgi
libexec/xymon/bb-csvinfo.cgi
libexec/xymon/bb-datepage.cgi
libexec/xymon/bb-eventlog.cgi
libexec/xymon/bb-findhost.cgi
libexec/xymon/bb-hist.cgi
libexec/xymon/bb-rep.cgi
libexec/xymon/bb-replog.cgi
libexec/xymon/bb-snapshot.cgi
libexec/xymon/bb-webpage
libexec/xymon/bbcmd
libexec/xymon/bbcombotest
libexec/xymon/bbdigest
libexec/xymon/bbgen
libexec/xymon/bbhostgrep
libexec/xymon/bbhostshow
libexec/xymon/bbmessage.cgi
libexec/xymon/bbproxy
libexec/xymon/bbtest-net
libexec/xymon/beastat
libexec/xymon/boilerplate.cgi
libexec/xymon/cgi-bin/bb-csvinfo.sh
libexec/xymon/cgi-bin/bb-datepage.sh
libexec/xymon/cgi-bin/bb-eventlog.sh
libexec/xymon/cgi-bin/bb-findhost.sh
libexec/xymon/cgi-bin/bb-hist.sh
libexec/xymon/cgi-bin/bb-histlog.sh
libexec/xymon/cgi-bin/bb-hostsvc.sh
libexec/xymon/cgi-bin/bb-rep.sh
libexec/xymon/cgi-bin/bb-replog.sh
libexec/xymon/cgi-bin/bb-snapshot.sh
libexec/xymon/cgi-bin/hobbit-certreport.sh
libexec/xymon/cgi-bin/hobbit-confreport-critical.sh
libexec/xymon/cgi-bin/hobbit-confreport.sh
libexec/xymon/cgi-bin/hobbit-ghosts.sh
libexec/xymon/cgi-bin/hobbit-hostgraphs.sh
libexec/xymon/cgi-bin/hobbit-hostlist.sh
libexec/xymon/cgi-bin/hobbit-nkview.sh
libexec/xymon/cgi-bin/hobbit-nongreen.sh
libexec/xymon/cgi-bin/hobbit-notifylog.sh
libexec/xymon/cgi-bin/hobbit-perfdata.sh
libexec/xymon/cgi-bin/hobbit-topchanges.sh
libexec/xymon/cgi-bin/hobbitcolumn.sh
libexec/xymon/cgi-bin/hobbitgraph.sh
libexec/xymon/cgi-secure/bb-ack.sh
libexec/xymon/cgi-secure/hobbit-ackinfo.sh
libexec/xymon/cgi-secure/hobbit-enadis.sh
libexec/xymon/cgi-secure/hobbit-nkedit.sh
libexec/xymon/cgi-secure/hobbit-useradm.sh
libexec/xymon/convertnk
libexec/xymon/ext/bbretest-net.sh
libexec/xymon/hobbit-ackinfo.cgi
libexec/xymon/hobbit-confreport.cgi
libexec/xymon/hobbit-enadis.cgi
libexec/xymon/hobbit-ghosts.cgi
libexec/xymon/hobbit-hostgraphs.cgi
libexec/xymon/hobbit-hostlist.cgi
libexec/xymon/hobbit-mailack
libexec/xymon/hobbit-nkedit.cgi
libexec/xymon/hobbit-nkview.cgi
libexec/xymon/hobbit-notifylog.cgi
libexec/xymon/hobbit-perfdata.cgi
libexec/xymon/hobbit-statusreport.cgi
libexec/xymon/hobbit-useradm.cgi
libexec/xymon/hobbit.sh
libexec/xymon/hobbitd
libexec/xymon/hobbitd_alert
libexec/xymon/hobbitd_capture
libexec/xymon/hobbitd_channel
libexec/xymon/hobbitd_client
libexec/xymon/hobbitd_filestore
libexec/xymon/hobbitd_history
libexec/xymon/hobbitd_hostdata
libexec/xymon/hobbitd_larrd
libexec/xymon/hobbitd_locator
libexec/xymon/hobbitd_rrd
libexec/xymon/hobbitd_sample
libexec/xymon/hobbitfetch
libexec/xymon/hobbitgraph.cgi
libexec/xymon/hobbitlaunch
libexec/xymon/hobbitping
libexec/xymon/hobbitreports.sh
libexec/xymon/hobbitsvc.cgi
libexec/xymon/moverrd.sh
libexec/xymon/rrdcachectl
libexec/xymon/trimhistory
man/man1/bb-ack.cgi.1
man/man1/bb-csvinfo.cgi.1
man/man1/bb-datepage.cgi.1
man/man1/bb-eventlog.cgi.1
man/man1/bb-findhost.cgi.1
man/man1/bb-hist.cgi.1
man/man1/bb-rep.cgi.1
man/man1/bb-replog.cgi.1
man/man1/bb-snapshot.cgi.1
man/man1/bb-webpage.cgi.1
man/man1/bb.1
man/man1/bbcmd.1
man/man1/bbcombotest.1
man/man1/bbdigest.1
man/man1/bbgen.1
man/man1/bbhostgrep.1
man/man1/bbhostshow.1
man/man1/bbretest-net.sh.1
man/man1/bbtest-net.1
man/man1/clientupdate.1
man/man1/hobbit-ackinfo.cgi.1
man/man1/hobbit-confreport.cgi.1
man/man1/hobbit-ghosts.cgi.1
man/man1/hobbit-hostgraphs.cgi.1
man/man1/hobbit-nkedit.cgi.1
man/man1/hobbit-nkview.cgi.1
man/man1/hobbit-statusreport.cgi.1
man/man1/hobbitgraph.cgi.1
man/man1/hobbitping.1
man/man1/hobbitsvc.cgi.1
man/man1/logfetch.1
man/man1/orcahobbit.1
man/man5/bb-hosts.5
man/man5/bb-services.5
man/man5/bbcombotest.cfg.5
man/man5/client-local.cfg.5
man/man5/clientlaunch.cfg.5
man/man5/hobbit-alerts.cfg.5
man/man5/hobbit-clients.cfg.5
man/man5/hobbit-nkview.cfg.5
man/man5/hobbitcgi.cfg.5
man/man5/hobbitclient.cfg.5
man/man5/hobbitgraph.cfg.5
man/man5/hobbitlaunch.cfg.5
man/man5/hobbitserver.cfg.5
man/man5/hobbitweb.5
man/man7/xymon.7
man/man8/bbmessage.cgi.8
man/man8/bbproxy.8
man/man8/hobbit-enadis.cgi.8
man/man8/hobbit-mailack.8
man/man8/hobbitd.8
man/man8/hobbitd_alert.8
man/man8/hobbitd_channel.8
man/man8/hobbitd_client.8
man/man8/hobbitd_filestore.8
man/man8/hobbitd_history.8
man/man8/hobbitd_hostdata.8
man/man8/hobbitd_rrd.8
man/man8/hobbitd_sample.8
man/man8/hobbitfetch.8
man/man8/hobbitlaunch.8
man/man8/msgcache.8
man/man8/trimhistory.8
share/examples/rc.d/xymon
share/examples/xymon/bb-hosts
share/examples/xymon/bb-services
share/examples/xymon/bbcombotest.cfg
share/examples/xymon/client-local.cfg
share/examples/xymon/columndoc.csv
share/examples/xymon/hobbit-alerts.cfg
share/examples/xymon/hobbit-apache.conf
share/examples/xymon/hobbit-clients.cfg
share/examples/xymon/hobbit-nkview.cfg
share/examples/xymon/hobbit-nkview.cfg.bak
share/examples/xymon/hobbitcgi.cfg
share/examples/xymon/hobbitgraph.cfg
share/examples/xymon/hobbitlaunch.cfg
share/examples/xymon/hobbitserver.cfg
share/examples/xymon/www/gifs/README
share/examples/xymon/www/gifs/arrow.gif
share/examples/xymon/www/gifs/bkg-blue.gif
share/examples/xymon/www/gifs/bkg-clear.gif
share/examples/xymon/www/gifs/bkg-green.gif
share/examples/xymon/www/gifs/bkg-purple.gif
share/examples/xymon/www/gifs/bkg-red.gif
share/examples/xymon/www/gifs/bkg-yellow.gif
share/examples/xymon/www/gifs/blue-ack.gif
share/examples/xymon/www/gifs/blue-recent.gif
share/examples/xymon/www/gifs/blue.gif
share/examples/xymon/www/gifs/clear-recent.gif
share/examples/xymon/www/gifs/clear.gif
share/examples/xymon/www/gifs/favicon-blue.ico
share/examples/xymon/www/gifs/favicon-clear.ico
share/examples/xymon/www/gifs/favicon-green.ico
share/examples/xymon/www/gifs/favicon-purple.ico
share/examples/xymon/www/gifs/favicon-red.ico
share/examples/xymon/www/gifs/favicon-unknown.ico
share/examples/xymon/www/gifs/favicon-yellow.ico
share/examples/xymon/www/gifs/green-recent.gif
share/examples/xymon/www/gifs/green.gif
share/examples/xymon/www/gifs/purple-ack.gif
share/examples/xymon/www/gifs/purple-recent.gif
share/examples/xymon/www/gifs/purple.gif
share/examples/xymon/www/gifs/red-ack.gif
share/examples/xymon/www/gifs/red-recent.gif
share/examples/xymon/www/gifs/red.gif
share/examples/xymon/www/gifs/unknown-recent.gif
share/examples/xymon/www/gifs/unknown.gif
share/examples/xymon/www/gifs/yellow-ack.gif
share/examples/xymon/www/gifs/yellow-recent.gif
share/examples/xymon/www/gifs/yellow.gif
share/examples/xymon/www/gifs/zoom.gif
share/examples/xymon/www/help/about.html
share/examples/xymon/www/help/bb-to-hobbit.html
share/examples/xymon/www/help/clonewarn.jpg
share/examples/xymon/www/help/configure.txt
share/examples/xymon/www/help/criticalsystems.html
share/examples/xymon/www/help/critview-detail-acked.jpg
share/examples/xymon/www/help/critview-detail-ackform.jpg
share/examples/xymon/www/help/critview-detail.jpg
share/examples/xymon/www/help/critview-disk.jpg
share/examples/xymon/www/help/critview-green.jpg
share/examples/xymon/www/help/editor-clonemaster.jpg
share/examples/xymon/www/help/editor-diskchanged.jpg
share/examples/xymon/www/help/editor-disksetup.jpg
share/examples/xymon/www/help/editor-main.jpg
share/examples/xymon/www/help/editor-makeclone.jpg
share/examples/xymon/www/help/editor-nohost.jpg
share/examples/xymon/www/help/editor-showclone.jpg
share/examples/xymon/www/help/hobbit-alerts.html
share/examples/xymon/www/help/hobbit-apacheconf.txt
share/examples/xymon/www/help/hobbit-bbhosts.png
share/examples/xymon/www/help/hobbit-clients.png
share/examples/xymon/www/help/hobbit-config.html
share/examples/xymon/www/help/hobbit-mrtg.html
share/examples/xymon/www/help/hobbit-tips.html
share/examples/xymon/www/help/hobbitmain.png
share/examples/xymon/www/help/hobbitprocs.png
share/examples/xymon/www/help/howtograph.html
share/examples/xymon/www/help/install.html
share/examples/xymon/www/help/known-issues.html
share/examples/xymon/www/help/mainview-acked.jpg
share/examples/xymon/www/help/mainview.jpg
share/examples/xymon/www/help/manpages/index.html
share/examples/xymon/www/help/manpages/man1/bb-ack.cgi.1.html
share/examples/xymon/www/help/manpages/man1/bb-csvinfo.cgi.1.html
share/examples/xymon/www/help/manpages/man1/bb-datepage.cgi.1.html
share/examples/xymon/www/help/manpages/man1/bb-eventlog.cgi.1.html
share/examples/xymon/www/help/manpages/man1/bb-findhost.cgi.1.html
share/examples/xymon/www/help/manpages/man1/bb-hist.cgi.1.html
share/examples/xymon/www/help/manpages/man1/bb-rep.cgi.1.html
share/examples/xymon/www/help/manpages/man1/bb-replog.cgi.1.html
share/examples/xymon/www/help/manpages/man1/bb-snapshot.cgi.1.html
share/examples/xymon/www/help/manpages/man1/bb-webpage.cgi.1.html
share/examples/xymon/www/help/manpages/man1/bb.1.html
share/examples/xymon/www/help/manpages/man1/bbcmd.1.html
share/examples/xymon/www/help/manpages/man1/bbcombotest.1.html
share/examples/xymon/www/help/manpages/man1/bbdigest.1.html
share/examples/xymon/www/help/manpages/man1/bbgen.1.html
share/examples/xymon/www/help/manpages/man1/bbhostgrep.1.html
share/examples/xymon/www/help/manpages/man1/bbhostshow.1.html
share/examples/xymon/www/help/manpages/man1/bbretest-net.sh.1.html
share/examples/xymon/www/help/manpages/man1/bbtest-net.1.html
share/examples/xymon/www/help/manpages/man1/clientupdate.1.html
share/examples/xymon/www/help/manpages/man1/hobbit-ackinfo.cgi.1.html
share/examples/xymon/www/help/manpages/man1/hobbit-confreport.cgi.1.html
share/examples/xymon/www/help/manpages/man1/hobbit-ghosts.cgi.1.html
share/examples/xymon/www/help/manpages/man1/hobbit-hostgraphs.cgi.1.html
share/examples/xymon/www/help/manpages/man1/hobbit-nkedit.cgi.1.html
share/examples/xymon/www/help/manpages/man1/hobbit-nkview.cgi.1.html
share/examples/xymon/www/help/manpages/man1/hobbit-statusreport.cgi.1.html
share/examples/xymon/www/help/manpages/man1/hobbitgraph.cgi.1.html
share/examples/xymon/www/help/manpages/man1/hobbitping.1.html
share/examples/xymon/www/help/manpages/man1/hobbitsvc.cgi.1.html
share/examples/xymon/www/help/manpages/man1/logfetch.1.html
share/examples/xymon/www/help/manpages/man1/orcahobbit.1.html
share/examples/xymon/www/help/manpages/man5/bb-hosts.5.html
share/examples/xymon/www/help/manpages/man5/bb-services.5.html
share/examples/xymon/www/help/manpages/man5/bbcombotest.cfg.5.html
share/examples/xymon/www/help/manpages/man5/client-local.cfg.5.html
share/examples/xymon/www/help/manpages/man5/clientlaunch.cfg.5.html
share/examples/xymon/www/help/manpages/man5/hobbit-alerts.cfg.5.html
share/examples/xymon/www/help/manpages/man5/hobbit-clients.cfg.5.html
share/examples/xymon/www/help/manpages/man5/hobbit-nkview.cfg.5.html
share/examples/xymon/www/help/manpages/man5/hobbitcgi.cfg.5.html
share/examples/xymon/www/help/manpages/man5/hobbitclient.cfg.5.html
share/examples/xymon/www/help/manpages/man5/hobbitgraph.cfg.5.html
share/examples/xymon/www/help/manpages/man5/hobbitlaunch.cfg.5.html
share/examples/xymon/www/help/manpages/man5/hobbitserver.cfg.5.html
share/examples/xymon/www/help/manpages/man5/hobbitweb.5.html
share/examples/xymon/www/help/manpages/man7/xymon.7.html
share/examples/xymon/www/help/manpages/man8/bbmessage.cgi.8.html
share/examples/xymon/www/help/manpages/man8/bbproxy.8.html
share/examples/xymon/www/help/manpages/man8/hobbit-enadis.cgi.8.html
share/examples/xymon/www/help/manpages/man8/hobbit-mailack.8.html
share/examples/xymon/www/help/manpages/man8/hobbitd.8.html
share/examples/xymon/www/help/manpages/man8/hobbitd_alert.8.html
share/examples/xymon/www/help/manpages/man8/hobbitd_channel.8.html
share/examples/xymon/www/help/manpages/man8/hobbitd_client.8.html
share/examples/xymon/www/help/manpages/man8/hobbitd_filestore.8.html
share/examples/xymon/www/help/manpages/man8/hobbitd_history.8.html
share/examples/xymon/www/help/manpages/man8/hobbitd_hostdata.8.html
share/examples/xymon/www/help/manpages/man8/hobbitd_rrd.8.html
share/examples/xymon/www/help/manpages/man8/hobbitd_sample.8.html
share/examples/xymon/www/help/manpages/man8/hobbitfetch.8.html
share/examples/xymon/www/help/manpages/man8/hobbitlaunch.8.html
share/examples/xymon/www/help/manpages/man8/msgcache.8.html
share/examples/xymon/www/help/manpages/man8/trimhistory.8.html
share/examples/xymon/www/help/stdview-detail-acked.jpg
share/examples/xymon/www/menu/README
share/examples/xymon/www/menu/menu.css
share/examples/xymon/www/menu/menu.js
share/examples/xymon/www/menu/menu_items.js
share/examples/xymon/www/menu/menu_tpl.js
share/xymon/web/acknowledge_footer
share/xymon/web/acknowledge_form
share/xymon/web/acknowledge_header
share/xymon/web/bb2_footer
share/xymon/web/bb2_header
share/xymon/web/bb_footer
share/xymon/web/bb_header
share/xymon/web/bbnk_footer
share/xymon/web/bbnk_header
share/xymon/web/bbrep_footer
share/xymon/web/bbrep_header
share/xymon/web/bbsnap2_footer
share/xymon/web/bbsnap2_header
share/xymon/web/bbsnap_footer
share/xymon/web/bbsnap_header
share/xymon/web/bbsnapnk_footer
share/xymon/web/bbsnapnk_header
share/xymon/web/columndoc_footer
share/xymon/web/columndoc_header
share/xymon/web/confreport_back
share/xymon/web/confreport_footer
share/xymon/web/confreport_front
share/xymon/web/confreport_header
share/xymon/web/event_footer
share/xymon/web/event_form
share/xymon/web/event_header
share/xymon/web/findhost_footer
share/xymon/web/findhost_form
share/xymon/web/findhost_header
share/xymon/web/ghosts_footer
share/xymon/web/ghosts_header
share/xymon/web/graphs_footer
share/xymon/web/graphs_header
share/xymon/web/hist_footer
share/xymon/web/hist_header
share/xymon/web/histlog_footer
share/xymon/web/histlog_header
share/xymon/web/hobbitnk_footer
share/xymon/web/hobbitnk_header
share/xymon/web/hostgraphs_footer
share/xymon/web/hostgraphs_form
share/xymon/web/hostgraphs_header
share/xymon/web/hostlist_footer
share/xymon/web/hostlist_form
share/xymon/web/hostlist_header
share/xymon/web/hostsvc_footer
share/xymon/web/hostsvc_header
share/xymon/web/info_footer
share/xymon/web/info_header
share/xymon/web/maint_footer
share/xymon/web/maint_form
share/xymon/web/maint_header
share/xymon/web/maintact_footer
share/xymon/web/maintact_header
share/xymon/web/nkack_form
share/xymon/web/nkedit_footer
share/xymon/web/nkedit_form
share/xymon/web/nkedit_header
share/xymon/web/notify_footer
share/xymon/web/notify_form
share/xymon/web/notify_header
share/xymon/web/perfdata_footer
share/xymon/web/perfdata_form
share/xymon/web/perfdata_header
share/xymon/web/replog_footer
share/xymon/web/replog_header
share/xymon/web/report_footer
share/xymon/web/report_form
share/xymon/web/report_form_daily
share/xymon/web/report_form_monthly
share/xymon/web/report_form_weekly
share/xymon/web/report_header
share/xymon/web/snapshot_footer
share/xymon/web/snapshot_form
share/xymon/web/snapshot_header
share/xymon/web/topchanges_footer
share/xymon/web/topchanges_form
share/xymon/web/topchanges_header
share/xymon/web/trends_form
share/xymon/web/useradm_footer
share/xymon/web/useradm_form
share/xymon/web/useradm_header
share/xymon/web/zoom.js
