# net/dnsmasq/PLIST
@comment $NetBSD: PLIST,v 1.1.1.1 2005/10/18 04:41:59 minskim Exp $
man/man8/dnsmasq.8
sbin/dnsmasq
share/examples/dnsmasq/dnsmasq.conf.example
share/examples/rc.d/dnsmasq
@dirrm share/examples/dnsmasq
