# net/dnsmasq/PLIST
@comment $NetBSD: PLIST,v 1.3 2012/01/31 12:39:22 roy Exp $
man/man8/dnsmasq.8
sbin/dnsmasq
share/examples/dnsmasq/dnsmasq.conf.example
share/examples/rc.d/dnsmasq
${PLIST.dbus}share/examples/dnsmasq/dnsmasq-dbus.conf
