# net/openslp/PLIST
@comment $NetBSD: PLIST,v 1.6 2005/05/02 20:34:03 reed Exp $
bin/slptool
include/slp.h
lib/libslp.la
sbin/slpd
share/doc/openslp/html/IntroductionToSLP/AgentInit.jpg
share/doc/openslp/html/IntroductionToSLP/ServiceReg.jpg
share/doc/openslp/html/IntroductionToSLP/ServiceRqst.jpg
share/doc/openslp/html/IntroductionToSLP/index.html
share/doc/openslp/html/ProgrammersGuide/Callbacks.html
share/doc/openslp/html/ProgrammersGuide/Divergence.html
share/doc/openslp/html/ProgrammersGuide/Examples.html
share/doc/openslp/html/ProgrammersGuide/SLPAttrCallback.html
share/doc/openslp/html/ProgrammersGuide/SLPClose.html
share/doc/openslp/html/ProgrammersGuide/SLPDelAttrs.html
share/doc/openslp/html/ProgrammersGuide/SLPDereg.html
share/doc/openslp/html/ProgrammersGuide/SLPError.html
share/doc/openslp/html/ProgrammersGuide/SLPEscape.html
share/doc/openslp/html/ProgrammersGuide/SLPFindAttrs.html
share/doc/openslp/html/ProgrammersGuide/SLPFindScopes.html
share/doc/openslp/html/ProgrammersGuide/SLPFindSrvTypes.html
share/doc/openslp/html/ProgrammersGuide/SLPFindSrvs.html
share/doc/openslp/html/ProgrammersGuide/SLPFree.html
share/doc/openslp/html/ProgrammersGuide/SLPGetProperty.html
share/doc/openslp/html/ProgrammersGuide/SLPGetRefreshInterval.html
share/doc/openslp/html/ProgrammersGuide/SLPOpen.html
share/doc/openslp/html/ProgrammersGuide/SLPParseSrvURL.html
share/doc/openslp/html/ProgrammersGuide/SLPReg.html
share/doc/openslp/html/ProgrammersGuide/SLPRegReport.html
share/doc/openslp/html/ProgrammersGuide/SLPSetProperty.html
share/doc/openslp/html/ProgrammersGuide/SLPSrvTypeCallback.html
share/doc/openslp/html/ProgrammersGuide/SLPSrvURLCallback.html
share/doc/openslp/html/ProgrammersGuide/SLPTypes.html
share/doc/openslp/html/ProgrammersGuide/SLPUnescape.html
share/doc/openslp/html/ProgrammersGuide/Security.html
share/doc/openslp/html/ProgrammersGuide/Syntax.html
share/doc/openslp/html/ProgrammersGuide/index.html
share/doc/openslp/html/ProgrammersGuide/openslp_security_whitepaper.html
share/doc/openslp/html/ProgrammersGuide/smalllogo.jpg
share/doc/openslp/html/UsersGuide/CommandLine.html
share/doc/openslp/html/UsersGuide/FileLocations.html
share/doc/openslp/html/UsersGuide/Installation.html
share/doc/openslp/html/UsersGuide/Optimization.html
share/doc/openslp/html/UsersGuide/Security.html
share/doc/openslp/html/UsersGuide/SlpConf.html
share/doc/openslp/html/UsersGuide/SlpReg.html
share/doc/openslp/html/UsersGuide/UserFAQ.html
share/doc/openslp/html/UsersGuide/WhenToRunSlpd.html
share/doc/openslp/html/UsersGuide/WhoShouldRead.html
share/doc/openslp/html/UsersGuide/index.html
share/doc/openslp/html/UsersGuide/smalllogo.jpg
share/doc/openslp/html/faq.html
share/doc/openslp/rfc/rfc1766.txt
share/doc/openslp/rfc/rfc2165.txt
share/doc/openslp/rfc/rfc2254.txt
share/doc/openslp/rfc/rfc2396.txt
share/doc/openslp/rfc/rfc2608.txt
share/doc/openslp/rfc/rfc2609.txt
share/doc/openslp/rfc/rfc2610.txt
share/doc/openslp/rfc/rfc2614.txt
share/doc/openslp/rfc/srvreg-integrity.txt
share/doc/openslp/rfc/threat_analysis_min_security.html
share/examples/openslp/slp.conf
share/examples/openslp/slp.reg
share/examples/openslp/slp.spi
share/examples/rc.d/slpd
@dirrm share/examples/openslp
@dirrm share/doc/openslp/rfc
@dirrm share/doc/openslp/html/UsersGuide
@dirrm share/doc/openslp/html/ProgrammersGuide
@dirrm share/doc/openslp/html/IntroductionToSLP
@dirrm share/doc/openslp/html
@dirrm share/doc/openslp
