# security/spiped/files/smf/manifest.xml
<?xml version="1.0"?>
<!DOCTYPE service_bundle SYSTEM "/usr/share/lib/xml/dtd/service_bundle.dtd.1">
<service_bundle type="manifest" name="@SMF_NAME@">
	<service name="@SMF_PREFIX@/@SMF_NAME@" type="service" version="1">
		<create_default_instance enabled="false"/>

		<dependency name="network" grouping="require_all" restart_on="refresh" type="service">
			<service_fmri value="svc:/milestone/network:default"/>
		</dependency>

		<dependency name="filesystem" grouping="require_all" restart_on="refresh" type="service">
			<service_fmri value="svc:/system/filesystem/local"/>
		</dependency>

		<exec_method type="method" name="start" exec="@PREFIX@/@SMF_METHOD_FILE.spiped@" timeout_seconds="5"/>
		<exec_method type="method" name="stop"	exec=":kill" timeout_seconds="60"/>

		<property_group name="application" type="application"></property_group>

		<stability value="Evolving"/>
		<template>
			<common_name>
				<loctext xml:lang="C">symmetrically encrypted and authenticated pipes</loctext>
			</common_name>
		</template>

	</service>
</service_bundle>
