# net/unfs3/PLIST
@comment $NetBSD: PLIST,v 1.3 2009/06/14 18:09:46 joerg Exp $
man/man7/tags.7
man/man8/unfsd.8
sbin/unfsd
share/examples/rc.d/unfsd
share/examples/unfs3/exports
