# audio/daapd/PLIST
@comment $NetBSD: PLIST,v 1.4 2009/06/14 17:32:04 joerg Exp $
bin/daapd
man/man8/daapd.8
share/examples/daapd/daapd.conf
share/examples/rc.d/daapd
