# audio/daapd/PLIST
@comment $NetBSD: PLIST,v 1.3 2005/05/02 20:33:57 reed Exp $
bin/daapd
man/man8/daapd.8
share/examples/daapd/daapd.conf
share/examples/rc.d/daapd
@dirrm share/examples/daapd
