# mail/dk-milter/PLIST
@comment $NetBSD: PLIST,v 1.1.1.1 2005/07/20 19:42:50 tv Exp $
man/man8/dk-milter.8
sbin/dk-milter
share/doc/dk-milter/FEATURES
share/doc/dk-milter/INSTALL
share/doc/dk-milter/KNOWNBUGS
share/doc/dk-milter/README
share/doc/dk-milter/RELEASE_NOTES
share/examples/rc.d/dk-milter
@dirrm share/doc/dk-milter
