# mail/dk-milter/PLIST
@comment $NetBSD: PLIST,v 1.3 2009/06/14 18:04:31 joerg Exp $
man/man8/dk-milter.8
sbin/dk-milter
sbin/gentxt.csh
share/doc/dk-milter/FEATURES
share/doc/dk-milter/INSTALL
share/doc/dk-milter/KNOWNBUGS
share/doc/dk-milter/README
share/doc/dk-milter/RELEASE_NOTES
share/examples/rc.d/dk-milter
