# mail/dk-milter/PLIST
@comment $NetBSD: PLIST,v 1.5 2015/03/15 15:28:57 mef Exp $
lib/libdk.a
man/man8/dk-milter.8
sbin/dk-milter
sbin/dk-milter-gentxt.csh
share/doc/dk-milter/FEATURES
share/doc/dk-milter/INSTALL
share/doc/dk-milter/KNOWNBUGS
share/doc/dk-milter/LICENSE
share/doc/dk-milter/README
share/doc/dk-milter/RELEASE_NOTES
share/doc/dk-milter/libdk/dk.html
share/doc/dk-milter/libdk/dk_alg_t.html
share/doc/dk-milter/libdk/dk_body.html
share/doc/dk-milter/libdk/dk_canon_t.html
share/doc/dk-milter/libdk/dk_close.html
share/doc/dk-milter/libdk/dk_eoh.html
share/doc/dk-milter/libdk/dk_eom.html
share/doc/dk-milter/libdk/dk_flags.html
share/doc/dk-milter/libdk/dk_free.html
share/doc/dk-milter/libdk/dk_geterror.html
share/doc/dk-milter/libdk/dk_gethdrs.html
share/doc/dk-milter/libdk/dk_getidentity.html
share/doc/dk-milter/libdk/dk_getsig.html
share/doc/dk-milter/libdk/dk_header.html
share/doc/dk-milter/libdk/dk_init.html
share/doc/dk-milter/libdk/dk_lib.html
share/doc/dk-milter/libdk/dk_options.html
share/doc/dk-milter/libdk/dk_query_t.html
share/doc/dk-milter/libdk/dk_sigkey_t.html
share/doc/dk-milter/libdk/dk_sign.html
share/doc/dk-milter/libdk/dk_stat.html
share/doc/dk-milter/libdk/dk_timeout.html
share/doc/dk-milter/libdk/dk_verify.html
share/doc/dk-milter/libdk/index.html
share/doc/dk-milter/libdk/overview.html
share/doc/dk-milter/libdk/rfc2822_mailbox_split.html
share/doc/dk-milter/rfc4870.txt
share/examples/rc.d/dk-milter
