# net/upclient/PLIST
@comment $NetBSD: PLIST,v 1.5 2005/05/02 20:34:04 reed Exp $
bin/upclient
share/doc/upclient/AUTHORS
share/doc/upclient/COPYING
share/doc/upclient/HISTORY
share/doc/upclient/README
share/doc/upclient/FAQ
share/doc/upclient/INSTALL
share/examples/rc.d/upclient
share/examples/upclient/upclient.conf
@dirrm share/examples/upclient
@dirrm share/doc/upclient
