# net/upclient/PLIST
@comment $NetBSD: PLIST,v 1.3 2003/07/02 18:22:42 jmmv Exp $
bin/upclient
etc/rc.d/upclient
share/doc/upclient/AUTHORS
share/doc/upclient/COPYING
share/doc/upclient/HISTORY
share/doc/upclient/README
share/doc/upclient/FAQ
share/doc/upclient/INSTALL
share/examples/upclient/upclient.conf
@dirrm share/examples/upclient
@dirrm share/doc/upclient
