# mail/qpopper/PLIST
@comment $NetBSD: PLIST,v 1.3 2005/05/02 20:34:01 reed Exp $
libexec/qpopper
man/man8/qpopper.8
man/man8/qpopauth.8
sbin/qpopauth
sbin/qpopper
share/doc/qpopper.License
share/examples/qpopper.config
share/examples/rc.d/qpopper
