# sysutils/ups-nut/PLIST
@comment $NetBSD: PLIST,v 1.4 2001/11/19 04:40:00 lukem Exp $
bin/aeg
bin/apcsmart
bin/belkin
bin/bestfort
bin/bestuferrups
bin/bestups
bin/engetron
bin/everups
bin/fentonups
bin/genericups
bin/ipt-anzen
bin/masterguard
bin/mge-ellipse
bin/mgeups
bin/multilink
bin/mustekups
bin/newapc
bin/optiups
bin/powercom
bin/sec
bin/sms
bin/toshiba1500
bin/ups-trust425+625
bin/upsc
bin/upscmd
bin/upsct
bin/upsct2
bin/upsdrvctl
bin/upseyeux
bin/upslog
bin/victronups
etc/nut/hosts.conf.sample
etc/nut/multimon.conf.sample
etc/nut/ups.conf.sample
etc/nut/upsd.conf.sample
etc/nut/upsd.users.sample
etc/nut/upsmon.conf.sample
etc/nut/upssched.conf.sample
etc/nut/upsset.conf.sample
etc/rc.d/upsd
man/man5/ups.conf.5
man/man5/upssched.conf.5
man/man8/apcsmart.8
man/man8/belkin.8
man/man8/bestups.8
man/man8/fentonups.8
man/man8/nutupsdrv.8
man/man8/powercom.8
man/man8/upsc.8
man/man8/upscmd.8
man/man8/upsct.8
man/man8/upsct2.8
man/man8/upsd.8
man/man8/upslog.8
man/man8/upsmon.8
man/man8/upssched.8
sbin/upsd
sbin/upsmon
sbin/upssched
share/doc/nut/Changes.trust
share/doc/nut/FAQ
share/doc/nut/README
share/doc/nut/access.txt
share/doc/nut/belkin.txt
share/doc/nut/big-servers.txt
share/doc/nut/cables/powerware.txt
share/doc/nut/cables/repotec.txt
share/doc/nut/cables/sms.txt
share/doc/nut/cables/victron.txt
share/doc/nut/chroot.txt
share/doc/nut/commands.txt
share/doc/nut/config-files.txt
share/doc/nut/configure.txt
share/doc/nut/data-room.txt
share/doc/nut/design.txt
share/doc/nut/driver-control.txt
share/doc/nut/dummycons.txt
share/doc/nut/dummyups.txt
share/doc/nut/generic-ups.txt
share/doc/nut/mge-ellipse.txt
share/doc/nut/model-arguments.txt
share/doc/nut/multilink.txt
share/doc/nut/new-modules.txt
share/doc/nut/pager.txt
share/doc/nut/powercom.txt
share/doc/nut/protocol.txt
share/doc/nut/shutdown.txt
share/doc/nut/tips.txt
share/doc/nut/todo.txt
share/doc/nut/ups-trust425+625.txt
share/doc/nut/upssched.txt
@dirrm share/doc/nut/cables
@dirrm share/doc/nut
