# databases/sqlrelay/PLIST
@comment $NetBSD: PLIST,v 1.8 2011/11/02 13:59:29 fhajny Exp $
bin/backupschema
bin/fields
bin/query
bin/sqlr-cachemanager
${PLIST.freetds}bin/sqlr-connection-freetds
${PLIST.mysql}bin/sqlr-connection-mysql
${PLIST.odbc}bin/sqlr-connection-odbc
${PLIST.pgsql}bin/sqlr-connection-postgresql
bin/sqlr-connection-router
${PLIST.sqlite}bin/sqlr-connection-sqlite
bin/sqlr-export
bin/sqlr-import
bin/sqlr-listener
bin/sqlr-scaler
bin/sqlr-start
bin/sqlr-status
bin/sqlr-stop
bin/sqlrclient-config
bin/sqlrclientwrapper-config
bin/sqlrsh
include/sqlrelay/private/bindvar.h
include/sqlrelay/private/column.h
include/sqlrelay/private/row.h
include/sqlrelay/private/sqlrconnection.h
include/sqlrelay/private/sqlrcursor.h
include/sqlrelay/private/sqlrdefines.h
include/sqlrelay/private/sqlrincludes.h
include/sqlrelay/sqlrclient.h
include/sqlrelay/sqlrclientwrapper.h
lib/libmysql3sqlrelay.la
lib/libmysql40sqlrelay.la
lib/libmysql41sqlrelay.la
lib/libmysql50sqlrelay.la
lib/libmysql51sqlrelay.la
lib/libpqsqlrelay.la
lib/libsqlrclient.la
lib/libsqlrclientwrapper.la
lib/libsqlrconnection.la
${PLIST.odbc}lib/libsqlrodbc.la
lib/libsqlrutil.la
lib/pkgconfig/sqlrelay-c++.pc
lib/pkgconfig/sqlrelay-c.pc
man/man1/fields.1
man/man1/query.1
man/man1/query.py.1
man/man1/sqlrsh.1
man/man8/sqlr-cachemanager.8
man/man8/sqlr-connection.8
man/man8/sqlr-ipclean.8
man/man8/sqlr-listener.8
man/man8/sqlr-scaler.8
man/man8/sqlr-start.8
man/man8/sqlr-stop.8
${PLIST.perl}${PERL5_SUB_INSTALLVENDORMAN3DIR}/SQLRelay::Connection.3
${PLIST.perl}${PERL5_SUB_INSTALLVENDORMAN3DIR}/SQLRelay::Cursor.3
${PLIST.perl}${PERL5_SUB_INSTALLVENDORARCH}/SQLRelay/Connection.pm
${PLIST.perl}${PERL5_SUB_INSTALLVENDORARCH}/SQLRelay/Cursor.pm
${PLIST.perl}${PERL5_SUB_INSTALLVENDORARCH}/auto/SQLRelay/Connection/Connection.so
${PLIST.perl}${PERL5_SUB_INSTALLVENDORARCH}/auto/SQLRelay/Cursor/Cursor.so
share/doc/sqlrelay/api/c++/sqlrclient.h.html
share/doc/sqlrelay/api/c/sqlrclientwrapper.h.html
share/doc/sqlrelay/api/java/allclasses-frame.html
share/doc/sqlrelay/api/java/allclasses-noframe.html
share/doc/sqlrelay/api/java/com/firstworks/sqlrelay/SQLRConnection.html
share/doc/sqlrelay/api/java/com/firstworks/sqlrelay/SQLRCursor.html
share/doc/sqlrelay/api/java/com/firstworks/sqlrelay/package-frame.html
share/doc/sqlrelay/api/java/com/firstworks/sqlrelay/package-summary.html
share/doc/sqlrelay/api/java/com/firstworks/sqlrelay/package-tree.html
share/doc/sqlrelay/api/java/constant-values.html
share/doc/sqlrelay/api/java/deprecated-list.html
share/doc/sqlrelay/api/java/help-doc.html
share/doc/sqlrelay/api/java/index-all.html
share/doc/sqlrelay/api/java/index.html
share/doc/sqlrelay/api/java/overview-tree.html
share/doc/sqlrelay/api/java/package-list
share/doc/sqlrelay/api/java/resources/inherit.gif
share/doc/sqlrelay/api/java/stylesheet.css
share/doc/sqlrelay/api/perl/Connection.html
share/doc/sqlrelay/api/perl/Cursor.html
share/doc/sqlrelay/api/perl/SQLRelay.html
share/doc/sqlrelay/api/perl/perl.html
share/doc/sqlrelay/api/php/doc.html
share/doc/sqlrelay/api/python/PySQLRClient.py.html
share/doc/sqlrelay/api/python/PySQLRDB.py.html
share/doc/sqlrelay/api/ruby/ruby.html
share/doc/sqlrelay/api/ruby/rubydbi.html
share/doc/sqlrelay/configuring.html
share/doc/sqlrelay/css/styles.css
share/doc/sqlrelay/dropin/mysql.html
share/doc/sqlrelay/dropin/postgresql.html
share/doc/sqlrelay/faq.html
share/doc/sqlrelay/gettingstarted/blackray.html
share/doc/sqlrelay/gettingstarted/db2.html
share/doc/sqlrelay/gettingstarted/db2/72.html
share/doc/sqlrelay/gettingstarted/db2/81.html
share/doc/sqlrelay/gettingstarted/db2/82.html
share/doc/sqlrelay/gettingstarted/db2/824.html
share/doc/sqlrelay/gettingstarted/db2/91.html
share/doc/sqlrelay/gettingstarted/firebird.html
share/doc/sqlrelay/gettingstarted/mssqlserver.html
share/doc/sqlrelay/gettingstarted/mysql.html
share/doc/sqlrelay/gettingstarted/oracle.html
share/doc/sqlrelay/gettingstarted/oracle/linux/fedoracore/1/10g.html
share/doc/sqlrelay/gettingstarted/oracle/linux/fedoracore/1/9iR2.html
share/doc/sqlrelay/gettingstarted/oracle/linux/fedoracore/2/10g.html
share/doc/sqlrelay/gettingstarted/oracle/linux/fedoracore/3/10g.html
share/doc/sqlrelay/gettingstarted/oracle/linux/fedoracore/5/10g.html
share/doc/sqlrelay/gettingstarted/oracle/linux/redhat/5/805.html
share/doc/sqlrelay/gettingstarted/oracle/linux/redhat/6/8i.html
share/doc/sqlrelay/gettingstarted/oracle/linux/redhat/7/8i.html
share/doc/sqlrelay/gettingstarted/oracle/linux/redhat/8/9iR2.html
share/doc/sqlrelay/gettingstarted/oracle/linux/redhat/9/9iR2.html
share/doc/sqlrelay/gettingstarted/oracle/linux/slackware/8/8i.html
share/doc/sqlrelay/gettingstarted/oracle/linux/slackware/8/9i.html
share/doc/sqlrelay/gettingstarted/oracle/linux/suse/9.0/9iR2.html
share/doc/sqlrelay/gettingstarted/postgresql.html
share/doc/sqlrelay/gettingstarted/sqlite.html
share/doc/sqlrelay/gettingstarted/sybase.html
share/doc/sqlrelay/gettingstarted/sybase/11.0.3.3.html
share/doc/sqlrelay/gettingstarted/sybase/11.9.2.html
share/doc/sqlrelay/gettingstarted/sybase/12.5.1.html
share/doc/sqlrelay/gettingstarted/sybase/12.5.2.html
share/doc/sqlrelay/gettingstarted/sybase/15.0.2.html
share/doc/sqlrelay/index.html
share/doc/sqlrelay/installing.html
share/doc/sqlrelay/introduction.html
share/doc/sqlrelay/introtowebdb.html
share/doc/sqlrelay/invisible.gif
share/doc/sqlrelay/oraclentier.html
share/doc/sqlrelay/programming/binds.html
share/doc/sqlrelay/programming/c++.html
share/doc/sqlrelay/programming/c.html
share/doc/sqlrelay/programming/java.html
share/doc/sqlrelay/programming/jdbc.html
share/doc/sqlrelay/programming/odbc.html
share/doc/sqlrelay/programming/perl.html
share/doc/sqlrelay/programming/perldbi.html
share/doc/sqlrelay/programming/php.html
share/doc/sqlrelay/programming/python.html
share/doc/sqlrelay/programming/pythondb.html
share/doc/sqlrelay/programming/ruby.html
share/doc/sqlrelay/programming/rubydbi.html
share/doc/sqlrelay/programming/tcl.html
share/doc/sqlrelay/programming/zope.html
share/doc/sqlrelay/programming/zope/instance.gif
share/doc/sqlrelay/programming/zope/management.gif
share/doc/sqlrelay/programming/zope/resultset.gif
share/doc/sqlrelay/programming/zope/sqlrelay.gif
share/doc/sqlrelay/programming/zope/test.gif
share/doc/sqlrelay/programming/zope/testzsqlmethod.gif
share/doc/sqlrelay/programming/zope/zsqlmethod.gif
share/doc/sqlrelay/programming/zope/zsqlmethodresultset.gif
share/doc/sqlrelay/running.html
share/doc/sqlrelay/tuning.html
share/doc/sqlrelay/windows.html
share/examples/rc.d/sqlrelay
share/examples/sqlrelay/sqlrelay
share/examples/sqlrelay/sqlrelay.conf.example
share/examples/sqlrelay/sqlrelay.dtd
