# databases/sqlrelay/files/smf/manifest.xml
<?xml version="1.0"?>
<!DOCTYPE service_bundle SYSTEM "/usr/share/lib/xml/dtd/service_bundle.dtd.1">
<service_bundle type="manifest" name="export">
  <service name="@SMF_PREFIX@/@SMF_NAME@" type="service" version="1">
    <dependency name="network" grouping="require_all" restart_on="error" type="service">
      <service_fmri value="svc:/milestone/network:default" />
    </dependency>
    <dependency name="filesystem" grouping="require_all" restart_on="error" type="service">
      <service_fmri value="svc:/system/filesystem/local" />
    </dependency>
    <method_context working_directory='@SQLR_DATA@'>
      <method_credential user='@SQLR_USER@' group='@SQLR_GROUP@' />
      <method_environment>
        <envvar name="PATH" value="@PREFIX@/sbin:@PREFIX@/bin:/usr/sbin:/usr/bin" />
      </method_environment>
    </method_context>
    <exec_method type="method" name="start" exec="@PREFIX@/bin/sqlr-start -id %i" timeout_seconds="3600" />
    <exec_method type="method" name="stop" exec=":kill" timeout_seconds="60" />
    <property_group name="startd" type="framework">
      <propval name="duration" type="astring" value="contract" />
      <propval name="ignore_error" type="astring" value="core,signal" />
    </property_group>
    <!--
         Instances should correspond to connections defined in @PKG_SYSCONFDIR@/sqlrelay.conf.
         For each connection you add there, create a new instance using svccfg, e.g.

             svccfg -s sqlrelay add mypgserver
    -->
    <instance name="example" enabled="false">
    </instance>
    <instance name="cachemanager" enabled="false">
      <exec_method type="method" name="start" exec="@PREFIX@/bin/sqlr-cachemanager" timeout_seconds="3600" />
      <exec_method type="method" name="stop" exec=":kill" timeout_seconds="3600" />
    </instance>
    <template>
      <common_name>
        <loctext xml:lang="C">SQL Relay</loctext>
      </common_name>
    </template>
  </service>
</service_bundle>
