# net/portmap/PLIST
@comment $NetBSD: PLIST,v 1.2 2009/06/14 18:09:38 joerg Exp $
sbin/pmap_dump
sbin/pmap_set
sbin/portmap
share/doc/portmap/README
share/examples/rc.d/portmap
