# sysutils/intel-microcode-netbsd/PLIST
@comment $NetBSD: PLIST,v 1.3 2014/02/19 19:00:14 drochner Exp $
libdata/firmware/cpu_x86_intel1/00000650-0
libdata/firmware/cpu_x86_intel1/00000650-1
libdata/firmware/cpu_x86_intel1/00000650-3
libdata/firmware/cpu_x86_intel1/00000651-0
libdata/firmware/cpu_x86_intel1/00000652-0
libdata/firmware/cpu_x86_intel1/00000652-1
libdata/firmware/cpu_x86_intel1/00000652-2
libdata/firmware/cpu_x86_intel1/00000653-0
libdata/firmware/cpu_x86_intel1/00000653-1
libdata/firmware/cpu_x86_intel1/00000653-2
libdata/firmware/cpu_x86_intel1/00000653-3
libdata/firmware/cpu_x86_intel1/00000660-0
libdata/firmware/cpu_x86_intel1/00000665-4
libdata/firmware/cpu_x86_intel1/0000066a-1
libdata/firmware/cpu_x86_intel1/0000066a-3
libdata/firmware/cpu_x86_intel1/0000066a-5
libdata/firmware/cpu_x86_intel1/0000066d-1
libdata/firmware/cpu_x86_intel1/0000066d-3
libdata/firmware/cpu_x86_intel1/0000066d-5
libdata/firmware/cpu_x86_intel1/00000671-2
libdata/firmware/cpu_x86_intel1/00000672-2
libdata/firmware/cpu_x86_intel1/00000673-2
libdata/firmware/cpu_x86_intel1/00000681-0
libdata/firmware/cpu_x86_intel1/00000681-2
libdata/firmware/cpu_x86_intel1/00000681-3
libdata/firmware/cpu_x86_intel1/00000681-4
libdata/firmware/cpu_x86_intel1/00000681-5
libdata/firmware/cpu_x86_intel1/00000683-3
libdata/firmware/cpu_x86_intel1/00000683-5
libdata/firmware/cpu_x86_intel1/00000686-0
libdata/firmware/cpu_x86_intel1/00000686-1
libdata/firmware/cpu_x86_intel1/00000686-2
libdata/firmware/cpu_x86_intel1/00000686-4
libdata/firmware/cpu_x86_intel1/00000686-7
libdata/firmware/cpu_x86_intel1/0000068a-4
libdata/firmware/cpu_x86_intel1/0000068a-5
libdata/firmware/cpu_x86_intel1/0000068a-7
libdata/firmware/cpu_x86_intel1/00000695-4
libdata/firmware/cpu_x86_intel1/00000695-5
libdata/firmware/cpu_x86_intel1/00000695-7
libdata/firmware/cpu_x86_intel1/000006a0-2
libdata/firmware/cpu_x86_intel1/000006a1-2
libdata/firmware/cpu_x86_intel1/000006b1-4
libdata/firmware/cpu_x86_intel1/000006b1-5
libdata/firmware/cpu_x86_intel1/000006b4-4
libdata/firmware/cpu_x86_intel1/000006b4-5
libdata/firmware/cpu_x86_intel1/000006d6-5
libdata/firmware/cpu_x86_intel1/000006e8-5
libdata/firmware/cpu_x86_intel1/000006ec-5
libdata/firmware/cpu_x86_intel1/000006ec-7
libdata/firmware/cpu_x86_intel1/000006f2-0
libdata/firmware/cpu_x86_intel1/000006f2-5
libdata/firmware/cpu_x86_intel1/000006f6-0
libdata/firmware/cpu_x86_intel1/000006f6-2
libdata/firmware/cpu_x86_intel1/000006f6-5
libdata/firmware/cpu_x86_intel1/000006f7-4
libdata/firmware/cpu_x86_intel1/000006f7-6
libdata/firmware/cpu_x86_intel1/000006fa-7
libdata/firmware/cpu_x86_intel1/000006fb-0
libdata/firmware/cpu_x86_intel1/000006fb-2
libdata/firmware/cpu_x86_intel1/000006fb-3
libdata/firmware/cpu_x86_intel1/000006fb-4
libdata/firmware/cpu_x86_intel1/000006fb-5
libdata/firmware/cpu_x86_intel1/000006fb-6
libdata/firmware/cpu_x86_intel1/000006fb-7
libdata/firmware/cpu_x86_intel1/000006fd-0
libdata/firmware/cpu_x86_intel1/000006fd-5
libdata/firmware/cpu_x86_intel1/000006fd-7
libdata/firmware/cpu_x86_intel1/00000f07-0
libdata/firmware/cpu_x86_intel1/00000f07-1
libdata/firmware/cpu_x86_intel1/00000f0a-0
libdata/firmware/cpu_x86_intel1/00000f0a-1
libdata/firmware/cpu_x86_intel1/00000f0a-2
libdata/firmware/cpu_x86_intel1/00000f12-2
libdata/firmware/cpu_x86_intel1/00000f24-1
libdata/firmware/cpu_x86_intel1/00000f24-2
libdata/firmware/cpu_x86_intel1/00000f24-4
libdata/firmware/cpu_x86_intel1/00000f25-0
libdata/firmware/cpu_x86_intel1/00000f25-1
libdata/firmware/cpu_x86_intel1/00000f25-2
libdata/firmware/cpu_x86_intel1/00000f25-4
libdata/firmware/cpu_x86_intel1/00000f26-1
libdata/firmware/cpu_x86_intel1/00000f27-1
libdata/firmware/cpu_x86_intel1/00000f27-2
libdata/firmware/cpu_x86_intel1/00000f27-3
libdata/firmware/cpu_x86_intel1/00000f29-1
libdata/firmware/cpu_x86_intel1/00000f29-2
libdata/firmware/cpu_x86_intel1/00000f29-3
libdata/firmware/cpu_x86_intel1/00000f32-0
libdata/firmware/cpu_x86_intel1/00000f32-2
libdata/firmware/cpu_x86_intel1/00000f32-3
libdata/firmware/cpu_x86_intel1/00000f33-0
libdata/firmware/cpu_x86_intel1/00000f33-2
libdata/firmware/cpu_x86_intel1/00000f33-3
libdata/firmware/cpu_x86_intel1/00000f34-0
libdata/firmware/cpu_x86_intel1/00000f34-2
libdata/firmware/cpu_x86_intel1/00000f34-3
libdata/firmware/cpu_x86_intel1/00000f34-4
libdata/firmware/cpu_x86_intel1/00000f41-0
libdata/firmware/cpu_x86_intel1/00000f41-1
libdata/firmware/cpu_x86_intel1/00000f41-2
libdata/firmware/cpu_x86_intel1/00000f41-3
libdata/firmware/cpu_x86_intel1/00000f41-4
libdata/firmware/cpu_x86_intel1/00000f41-5
libdata/firmware/cpu_x86_intel1/00000f41-7
libdata/firmware/cpu_x86_intel1/00000f43-0
libdata/firmware/cpu_x86_intel1/00000f43-2
libdata/firmware/cpu_x86_intel1/00000f43-3
libdata/firmware/cpu_x86_intel1/00000f43-4
libdata/firmware/cpu_x86_intel1/00000f43-7
libdata/firmware/cpu_x86_intel1/00000f44-0
libdata/firmware/cpu_x86_intel1/00000f44-2
libdata/firmware/cpu_x86_intel1/00000f44-3
libdata/firmware/cpu_x86_intel1/00000f44-4
libdata/firmware/cpu_x86_intel1/00000f44-7
libdata/firmware/cpu_x86_intel1/00000f47-0
libdata/firmware/cpu_x86_intel1/00000f47-2
libdata/firmware/cpu_x86_intel1/00000f47-3
libdata/firmware/cpu_x86_intel1/00000f47-4
libdata/firmware/cpu_x86_intel1/00000f47-7
libdata/firmware/cpu_x86_intel1/00000f48-0
libdata/firmware/cpu_x86_intel1/00000f48-1
libdata/firmware/cpu_x86_intel1/00000f48-2
libdata/firmware/cpu_x86_intel1/00000f48-3
libdata/firmware/cpu_x86_intel1/00000f48-4
libdata/firmware/cpu_x86_intel1/00000f48-6
libdata/firmware/cpu_x86_intel1/00000f49-0
libdata/firmware/cpu_x86_intel1/00000f49-2
libdata/firmware/cpu_x86_intel1/00000f49-3
libdata/firmware/cpu_x86_intel1/00000f49-4
libdata/firmware/cpu_x86_intel1/00000f49-5
libdata/firmware/cpu_x86_intel1/00000f49-7
libdata/firmware/cpu_x86_intel1/00000f4a-0
libdata/firmware/cpu_x86_intel1/00000f4a-2
libdata/firmware/cpu_x86_intel1/00000f4a-3
libdata/firmware/cpu_x86_intel1/00000f4a-4
libdata/firmware/cpu_x86_intel1/00000f4a-6
libdata/firmware/cpu_x86_intel1/00000f62-2
libdata/firmware/cpu_x86_intel1/00000f64-0
libdata/firmware/cpu_x86_intel1/00000f64-2
libdata/firmware/cpu_x86_intel1/00000f64-4
libdata/firmware/cpu_x86_intel1/00000f64-5
libdata/firmware/cpu_x86_intel1/00000f65-0
libdata/firmware/cpu_x86_intel1/00000f68-1
libdata/firmware/cpu_x86_intel1/00000f68-5
libdata/firmware/cpu_x86_intel1/00010661-0
libdata/firmware/cpu_x86_intel1/00010661-1
libdata/firmware/cpu_x86_intel1/00010661-7
libdata/firmware/cpu_x86_intel1/00010676-0
libdata/firmware/cpu_x86_intel1/00010676-2
libdata/firmware/cpu_x86_intel1/00010676-4
libdata/firmware/cpu_x86_intel1/00010676-6
libdata/firmware/cpu_x86_intel1/00010676-7
libdata/firmware/cpu_x86_intel1/00010677-4
libdata/firmware/cpu_x86_intel1/0001067a-0
libdata/firmware/cpu_x86_intel1/0001067a-2
libdata/firmware/cpu_x86_intel1/0001067a-4
libdata/firmware/cpu_x86_intel1/0001067a-5
libdata/firmware/cpu_x86_intel1/0001067a-6
libdata/firmware/cpu_x86_intel1/0001067a-7
libdata/firmware/cpu_x86_intel1/000106a4-0
libdata/firmware/cpu_x86_intel1/000106a4-1
libdata/firmware/cpu_x86_intel1/000106a5-0
libdata/firmware/cpu_x86_intel1/000106a5-1
libdata/firmware/cpu_x86_intel1/000106c2-0
libdata/firmware/cpu_x86_intel1/000106c2-2
libdata/firmware/cpu_x86_intel1/000106c2-3
libdata/firmware/cpu_x86_intel1/000106ca-0
libdata/firmware/cpu_x86_intel1/000106ca-2
libdata/firmware/cpu_x86_intel1/000106ca-3
libdata/firmware/cpu_x86_intel1/000106ca-4
libdata/firmware/cpu_x86_intel1/000106d1-3
libdata/firmware/cpu_x86_intel1/000106e4-0
libdata/firmware/cpu_x86_intel1/000106e4-3
libdata/firmware/cpu_x86_intel1/000106e5-0
libdata/firmware/cpu_x86_intel1/000106e5-1
libdata/firmware/cpu_x86_intel1/000106e5-4
libdata/firmware/cpu_x86_intel1/00020652-1
libdata/firmware/cpu_x86_intel1/00020652-4
libdata/firmware/cpu_x86_intel1/00020655-1
libdata/firmware/cpu_x86_intel1/00020655-4
libdata/firmware/cpu_x86_intel1/00020655-7
libdata/firmware/cpu_x86_intel1/00020661-0
libdata/firmware/cpu_x86_intel1/00020661-1
libdata/firmware/cpu_x86_intel1/000206a7-1
libdata/firmware/cpu_x86_intel1/000206a7-4
libdata/firmware/cpu_x86_intel1/000206d6-0
libdata/firmware/cpu_x86_intel1/000206d6-2
libdata/firmware/cpu_x86_intel1/000206d6-3
libdata/firmware/cpu_x86_intel1/000206d6-5
libdata/firmware/cpu_x86_intel1/000206d6-6
libdata/firmware/cpu_x86_intel1/000206d7-0
libdata/firmware/cpu_x86_intel1/000206d7-2
libdata/firmware/cpu_x86_intel1/000206d7-3
libdata/firmware/cpu_x86_intel1/000206d7-5
libdata/firmware/cpu_x86_intel1/000206d7-6
libdata/firmware/cpu_x86_intel1/000206f2-0
libdata/firmware/cpu_x86_intel1/000206f2-2
libdata/firmware/cpu_x86_intel1/000306a9-1
libdata/firmware/cpu_x86_intel1/000306a9-4
libdata/firmware/cpu_x86_intel1/000306c3-1
libdata/firmware/cpu_x86_intel1/000306c3-4
libdata/firmware/cpu_x86_intel1/000306c3-5
libdata/firmware/cpu_x86_intel1/000306e4-0
libdata/firmware/cpu_x86_intel1/000306e4-2
libdata/firmware/cpu_x86_intel1/000306e4-3
libdata/firmware/cpu_x86_intel1/000306e4-5
libdata/firmware/cpu_x86_intel1/000306e4-6
libdata/firmware/cpu_x86_intel1/000306e4-7
libdata/firmware/cpu_x86_intel1/000306e6-0
libdata/firmware/cpu_x86_intel1/000306e6-2
libdata/firmware/cpu_x86_intel1/000306e6-3
libdata/firmware/cpu_x86_intel1/000306e6-5
libdata/firmware/cpu_x86_intel1/000306e6-6
libdata/firmware/cpu_x86_intel1/000306e6-7
libdata/firmware/cpu_x86_intel1/00040651-1
libdata/firmware/cpu_x86_intel1/00040651-4
libdata/firmware/cpu_x86_intel1/00040651-5
libdata/firmware/cpu_x86_intel1/00040651-6
libdata/firmware/cpu_x86_intel1/00040661-1
libdata/firmware/cpu_x86_intel1/00040661-4
libdata/firmware/cpu_x86_intel1/00040661-5
share/examples/rc.d/intel-microcode
