# net/istgt/PLIST
@comment $NetBSD: PLIST,v 1.1.1.1 2010/02/02 16:09:51 sborrill Exp $
bin/istgt
bin/istgtcontrol
man/man1/istgt.1
share/examples/istgt/auth.conf
share/examples/istgt/istgt.conf
share/examples/istgt/istgt.large.conf
share/examples/istgt/istgtcontrol.conf
share/examples/rc.d/istgt
