# mail/amavisd-new/PLIST
@comment $NetBSD: PLIST,v 1.4 2011/10/25 20:33:07 spz Exp $
${PLIST.milter}sbin/amavis
${PLIST.milter}sbin/amavis-milter
sbin/amavisd
sbin/amavisd-agent
sbin/amavisd-nanny
sbin/amavisd-release
share/doc/amavisd-new/AAAREADME.first
share/doc/amavisd-new/INSTALL
share/doc/amavisd-new/LDAP.schema
share/doc/amavisd-new/LICENSE
share/doc/amavisd-new/MANIFEST
share/doc/amavisd-new/README.banned
share/doc/amavisd-new/README.chroot
share/doc/amavisd-new/README.contributed
share/doc/amavisd-new/README.courier
share/doc/amavisd-new/README.courier-old
share/doc/amavisd-new/README.customize
share/doc/amavisd-new/README.exim_v3
share/doc/amavisd-new/README.exim_v3_app
share/doc/amavisd-new/README.exim_v4
share/doc/amavisd-new/README.exim_v4_app
share/doc/amavisd-new/README.exim_v4_app2
share/doc/amavisd-new/README.ldap
share/doc/amavisd-new/README.lookups
share/doc/amavisd-new/README.milter
share/doc/amavisd-new/README.old.scanners
share/doc/amavisd-new/README.performance
share/doc/amavisd-new/README.policy-on-notifications
share/doc/amavisd-new/README.postfix
share/doc/amavisd-new/README.postfix.html
share/doc/amavisd-new/README.protocol
share/doc/amavisd-new/README.sendmail
share/doc/amavisd-new/README.sendmail-dual
share/doc/amavisd-new/README.sendmail-dual.old
share/doc/amavisd-new/README.sql
share/doc/amavisd-new/README.sql-mysql
share/doc/amavisd-new/README.sql-pg
share/doc/amavisd-new/RELEASE_NOTES
share/doc/amavisd-new/amavisd-new-docs.html
share/examples/amavisd-new/amavisd.conf-default
share/examples/amavisd-new/amavisd.conf-minimal
share/examples/rc.d/amavisd
${PLIST.milter}share/examples/rc.d/amavismilter
