# net/isc-dhclient/PLIST
@comment $NetBSD: PLIST,v 1.1.1.1 2007/12/12 20:03:48 adrianp Exp $
man/man5/dhclient.conf.5
man/man5/dhclient.leases.5
man/man8/dhclient-script.8
man/man8/dhclient.8
sbin/dhclient
sbin/dhclient-script
share/examples/rc.d/isc_dhclient
share/examples/isc-dhcp-client/dhclient.conf
@dirrm share/examples/isc-dhcp-client
