# sysutils/xenstoretools/PLIST
@comment $NetBSD: PLIST,v 1.2 2012/02/24 15:26:28 drochner Exp $
bin/xenstore
bin/xenstore-chmod
bin/xenstore-control
bin/xenstore-exists
bin/xenstore-list
bin/xenstore-ls
bin/xenstore-read
bin/xenstore-rm
bin/xenstore-watch
bin/xenstore-write
lib/libxenstore.so
lib/libxenstore.so.3.0
lib/libxenstore.so.3.0.0
share/examples/rc.d/xendomname
