# chat/atheme/PLIST
@comment $NetBSD: PLIST,v 1.2 2009/06/14 17:37:42 joerg Exp $
bin/atheme-services
lib/atheme/modules/backend/flatfile.so
lib/atheme/modules/chanserv/akick.so
lib/atheme/modules/chanserv/ban.so
lib/atheme/modules/chanserv/clear.so
lib/atheme/modules/chanserv/clear_bans.so
lib/atheme/modules/chanserv/clear_users.so
lib/atheme/modules/chanserv/close.so
lib/atheme/modules/chanserv/count.so
lib/atheme/modules/chanserv/drop.so
lib/atheme/modules/chanserv/fflags.so
lib/atheme/modules/chanserv/flags.so
lib/atheme/modules/chanserv/ftransfer.so
lib/atheme/modules/chanserv/getkey.so
lib/atheme/modules/chanserv/halfop.so
lib/atheme/modules/chanserv/help.so
lib/atheme/modules/chanserv/hold.so
lib/atheme/modules/chanserv/info.so
lib/atheme/modules/chanserv/invite.so
lib/atheme/modules/chanserv/kick.so
lib/atheme/modules/chanserv/list.so
lib/atheme/modules/chanserv/main.so
lib/atheme/modules/chanserv/mark.so
lib/atheme/modules/chanserv/op.so
lib/atheme/modules/chanserv/owner.so
lib/atheme/modules/chanserv/protect.so
lib/atheme/modules/chanserv/quiet.so
lib/atheme/modules/chanserv/recover.so
lib/atheme/modules/chanserv/register.so
lib/atheme/modules/chanserv/set.so
lib/atheme/modules/chanserv/set_private.so
lib/atheme/modules/chanserv/status.so
lib/atheme/modules/chanserv/taxonomy.so
lib/atheme/modules/chanserv/template.so
lib/atheme/modules/chanserv/topic.so
lib/atheme/modules/chanserv/unban_self.so
lib/atheme/modules/chanserv/version.so
lib/atheme/modules/chanserv/voice.so
lib/atheme/modules/chanserv/why.so
lib/atheme/modules/chanserv/xop.so
lib/atheme/modules/crypto/ircservices.so
lib/atheme/modules/crypto/posix.so
lib/atheme/modules/crypto/sorservices.so
lib/atheme/modules/gameserv/dice.so
lib/atheme/modules/gameserv/eightball.so
lib/atheme/modules/gameserv/help.so
lib/atheme/modules/gameserv/main.so
lib/atheme/modules/gameserv/namegen.so
lib/atheme/modules/gameserv/rps.so
lib/atheme/modules/global/main.so
lib/atheme/modules/memoserv/delete.so
lib/atheme/modules/memoserv/forward.so
lib/atheme/modules/memoserv/help.so
lib/atheme/modules/memoserv/ignore.so
lib/atheme/modules/memoserv/list.so
lib/atheme/modules/memoserv/main.so
lib/atheme/modules/memoserv/read.so
lib/atheme/modules/memoserv/send.so
lib/atheme/modules/memoserv/sendops.so
lib/atheme/modules/misc/httpd.so
lib/atheme/modules/nickserv/access.so
lib/atheme/modules/nickserv/drop.so
lib/atheme/modules/nickserv/enforce.so
lib/atheme/modules/nickserv/freeze.so
lib/atheme/modules/nickserv/ghost.so
lib/atheme/modules/nickserv/group.so
lib/atheme/modules/nickserv/help.so
lib/atheme/modules/nickserv/hold.so
lib/atheme/modules/nickserv/identify.so
lib/atheme/modules/nickserv/info.so
lib/atheme/modules/nickserv/list.so
lib/atheme/modules/nickserv/listchans.so
lib/atheme/modules/nickserv/listmail.so
lib/atheme/modules/nickserv/login.so
lib/atheme/modules/nickserv/logout.so
lib/atheme/modules/nickserv/main.so
lib/atheme/modules/nickserv/mark.so
lib/atheme/modules/nickserv/register.so
lib/atheme/modules/nickserv/resetpass.so
lib/atheme/modules/nickserv/return.so
lib/atheme/modules/nickserv/sendpass.so
lib/atheme/modules/nickserv/set.so
lib/atheme/modules/nickserv/set_accountname.so
lib/atheme/modules/nickserv/set_private.so
lib/atheme/modules/nickserv/set_privmsg.so
lib/atheme/modules/nickserv/setpass.so
lib/atheme/modules/nickserv/status.so
lib/atheme/modules/nickserv/subscribe.so
lib/atheme/modules/nickserv/taxonomy.so
lib/atheme/modules/nickserv/verify.so
lib/atheme/modules/nickserv/vhost.so
lib/atheme/modules/operserv/akill.so
lib/atheme/modules/operserv/clearchan.so
lib/atheme/modules/operserv/clones.so
lib/atheme/modules/operserv/compare.so
lib/atheme/modules/operserv/help.so
lib/atheme/modules/operserv/identify.so
lib/atheme/modules/operserv/ignore.so
lib/atheme/modules/operserv/inject.so
lib/atheme/modules/operserv/jupe.so
lib/atheme/modules/operserv/main.so
lib/atheme/modules/operserv/mode.so
lib/atheme/modules/operserv/modinspect.so
lib/atheme/modules/operserv/modlist.so
lib/atheme/modules/operserv/modload.so
lib/atheme/modules/operserv/modrestart.so
lib/atheme/modules/operserv/modunload.so
lib/atheme/modules/operserv/noop.so
lib/atheme/modules/operserv/rakill.so
lib/atheme/modules/operserv/raw.so
lib/atheme/modules/operserv/rehash.so
lib/atheme/modules/operserv/restart.so
lib/atheme/modules/operserv/rmatch.so
lib/atheme/modules/operserv/rnc.so
lib/atheme/modules/operserv/rwatch.so
lib/atheme/modules/operserv/shutdown.so
lib/atheme/modules/operserv/soper.so
lib/atheme/modules/operserv/specs.so
lib/atheme/modules/operserv/update.so
lib/atheme/modules/operserv/uptime.so
lib/atheme/modules/protocol/asuka.so
lib/atheme/modules/protocol/bahamut.so
lib/atheme/modules/protocol/bircd.so
lib/atheme/modules/protocol/charybdis.so
lib/atheme/modules/protocol/dreamforge.so
lib/atheme/modules/protocol/hybrid.so
lib/atheme/modules/protocol/hyperion.so
lib/atheme/modules/protocol/inspircd11.so
lib/atheme/modules/protocol/inspircd12.so
lib/atheme/modules/protocol/ircnet.so
lib/atheme/modules/protocol/nefarious.so
lib/atheme/modules/protocol/officeirc.so
lib/atheme/modules/protocol/plexus.so
lib/atheme/modules/protocol/ptlink.so
lib/atheme/modules/protocol/ratbox.so
lib/atheme/modules/protocol/shadowircd.so
lib/atheme/modules/protocol/solidircd.so
lib/atheme/modules/protocol/sorcery.so
lib/atheme/modules/protocol/ultimate3.so
lib/atheme/modules/protocol/undernet.so
lib/atheme/modules/protocol/unreal.so
lib/atheme/modules/saslserv/crypt.so
lib/atheme/modules/saslserv/dh-blowfish.so
lib/atheme/modules/saslserv/main.so
lib/atheme/modules/saslserv/plain.so
lib/atheme/modules/xmlrpc/account.so
lib/atheme/modules/xmlrpc/channel.so
lib/atheme/modules/xmlrpc/main.so
lib/atheme/modules/xmlrpc/memo.so
share/atheme/help/cservice/akick
share/atheme/help/cservice/ban
share/atheme/help/cservice/clear_bans
share/atheme/help/cservice/clear_users
share/atheme/help/cservice/close
share/atheme/help/cservice/count
share/atheme/help/cservice/drop
share/atheme/help/cservice/fflags
share/atheme/help/cservice/flags
share/atheme/help/cservice/forcexop
share/atheme/help/cservice/ftransfer
share/atheme/help/cservice/getkey
share/atheme/help/cservice/halfop
share/atheme/help/cservice/hold
share/atheme/help/cservice/info
share/atheme/help/cservice/invite
share/atheme/help/cservice/kick
share/atheme/help/cservice/kickban
share/atheme/help/cservice/list
share/atheme/help/cservice/mark
share/atheme/help/cservice/op_voice
share/atheme/help/cservice/quiet
share/atheme/help/cservice/recover
share/atheme/help/cservice/register
share/atheme/help/cservice/set_email
share/atheme/help/cservice/set_entrymsg
share/atheme/help/cservice/set_fantasy
share/atheme/help/cservice/set_founder
share/atheme/help/cservice/set_guard
share/atheme/help/cservice/set_keeptopic
share/atheme/help/cservice/set_mlock
share/atheme/help/cservice/set_private
share/atheme/help/cservice/set_property
share/atheme/help/cservice/set_secure
share/atheme/help/cservice/set_staffonly
share/atheme/help/cservice/set_topiclock
share/atheme/help/cservice/set_url
share/atheme/help/cservice/set_verbose
share/atheme/help/cservice/status
share/atheme/help/cservice/taxonomy
share/atheme/help/cservice/template
share/atheme/help/cservice/topic
share/atheme/help/cservice/topicappend
share/atheme/help/cservice/topicprepend
share/atheme/help/cservice/unban
share/atheme/help/cservice/unban_self
share/atheme/help/cservice/unquiet
share/atheme/help/cservice/userinfo
share/atheme/help/cservice/why
share/atheme/help/cservice/xop
share/atheme/help/gameserv/eightball
share/atheme/help/gameserv/namegen
share/atheme/help/gameserv/roll
share/atheme/help/gservice/global
share/atheme/help/help
share/atheme/help/memoserv/delete
share/atheme/help/memoserv/forward
share/atheme/help/memoserv/ignore
share/atheme/help/memoserv/list
share/atheme/help/memoserv/read
share/atheme/help/memoserv/send
share/atheme/help/memoserv/sendops
share/atheme/help/nickserv/acc
share/atheme/help/nickserv/access
share/atheme/help/nickserv/drop
share/atheme/help/nickserv/freeze
share/atheme/help/nickserv/fungroup
share/atheme/help/nickserv/fverify
share/atheme/help/nickserv/ghost
share/atheme/help/nickserv/group
share/atheme/help/nickserv/hold
share/atheme/help/nickserv/identify
share/atheme/help/nickserv/info
share/atheme/help/nickserv/list
share/atheme/help/nickserv/listchans
share/atheme/help/nickserv/listmail
share/atheme/help/nickserv/listvhost
share/atheme/help/nickserv/login
share/atheme/help/nickserv/logout
share/atheme/help/nickserv/mark
share/atheme/help/nickserv/register
share/atheme/help/nickserv/release
share/atheme/help/nickserv/resetpass
share/atheme/help/nickserv/return
share/atheme/help/nickserv/sendpass
share/atheme/help/nickserv/set_email
share/atheme/help/nickserv/set_emailmemos
share/atheme/help/nickserv/set_enforce
share/atheme/help/nickserv/set_hidemail
share/atheme/help/nickserv/set_neverop
share/atheme/help/nickserv/set_nomemo
share/atheme/help/nickserv/set_noop
share/atheme/help/nickserv/set_password
share/atheme/help/nickserv/set_private
share/atheme/help/nickserv/set_privmsg
share/atheme/help/nickserv/set_property
share/atheme/help/nickserv/setpass
share/atheme/help/nickserv/status
share/atheme/help/nickserv/taxonomy
share/atheme/help/nickserv/ungroup
share/atheme/help/nickserv/verify
share/atheme/help/nickserv/vhost
share/atheme/help/oservice/akill
share/atheme/help/oservice/clearchan
share/atheme/help/oservice/clones
share/atheme/help/oservice/compare
share/atheme/help/oservice/identify
share/atheme/help/oservice/ignore
share/atheme/help/oservice/inject
share/atheme/help/oservice/jupe
share/atheme/help/oservice/mode
share/atheme/help/oservice/modinspect
share/atheme/help/oservice/modlist
share/atheme/help/oservice/modload
share/atheme/help/oservice/modrestart
share/atheme/help/oservice/modunload
share/atheme/help/oservice/noop
share/atheme/help/oservice/rakill
share/atheme/help/oservice/raw
share/atheme/help/oservice/rehash
share/atheme/help/oservice/restart
share/atheme/help/oservice/rmatch
share/atheme/help/oservice/rnc
share/atheme/help/oservice/rwatch
share/atheme/help/oservice/shutdown
share/atheme/help/oservice/soper
share/atheme/help/oservice/specs
share/atheme/help/oservice/update
share/atheme/help/oservice/uptime
share/atheme/help/saslserv/pkey
share/atheme/help/userserv/drop
share/atheme/help/userserv/freeze
share/atheme/help/userserv/fverify
share/atheme/help/userserv/hold
share/atheme/help/userserv/info
share/atheme/help/userserv/list
share/atheme/help/userserv/listchans
share/atheme/help/userserv/listmail
share/atheme/help/userserv/listvhost
share/atheme/help/userserv/login
share/atheme/help/userserv/logout
share/atheme/help/userserv/mark
share/atheme/help/userserv/register
share/atheme/help/userserv/resetpass
share/atheme/help/userserv/return
share/atheme/help/userserv/sendpass
share/atheme/help/userserv/set_email
share/atheme/help/userserv/set_emailmemos
share/atheme/help/userserv/set_hidemail
share/atheme/help/userserv/set_neverop
share/atheme/help/userserv/set_noop
share/atheme/help/userserv/set_password
share/atheme/help/userserv/set_private
share/atheme/help/userserv/set_privmsg
share/atheme/help/userserv/set_property
share/atheme/help/userserv/setpass
share/atheme/help/userserv/status
share/atheme/help/userserv/taxonomy
share/atheme/help/userserv/verify
share/atheme/help/userserv/vhost
share/doc/atheme/INSTALL
share/doc/atheme/IRCD
share/doc/atheme/LICENSE
share/doc/atheme/PCOMMANDS
share/doc/atheme/PRIVILEGES
share/doc/atheme/README
share/doc/atheme/RELEASE
share/doc/atheme/SASL
share/doc/atheme/SQL
share/doc/atheme/TRANSLATION
share/doc/atheme/VERSION_LETTERS
share/doc/atheme/XMLRPC
share/doc/atheme/technical/HOOKS
share/doc/atheme/technical/MODES
share/doc/atheme/technical/XMLRPCLIB
share/examples/atheme/atheme.conf.example
share/examples/atheme/atheme.conf.userserv-example
share/examples/atheme/atheme.cron.example
share/examples/atheme/atheme.motd.example
share/examples/rc.d/atheme
