# www/nostromo/PLIST
@comment $NetBSD: PLIST,v 1.2 2019/10/20 20:02:13 ast Exp $
man/man8/nhttpd.8
sbin/crypt
sbin/nhttpd
share/examples/nostromo/conf/mimes
share/examples/nostromo/conf/nhttpd.conf-dist
share/examples/nostromo/htdocs/cgi-bin/printenv
share/examples/nostromo/htdocs/index.html
share/examples/nostromo/htdocs/nostromo.gif
share/examples/nostromo/icons/dir.gif
share/examples/nostromo/icons/file.gif
share/examples/rc.d/nostromo
