# audio/ubs/PLIST
@comment $NetBSD: PLIST,v 1.2 2005/05/02 20:33:57 reed Exp $
bin/ubs
bin/ubs-event
bin/ubs-play
bin/ubs-sched
include/defaults.h
include/ubs-event.h
include/ubs-main.h
include/ubs-play.h
include/ubs-sched.h
include/ubs.h
lib/libubs.a
share/examples/rc.d/ubs
share/examples/ubs/help.conf
share/examples/ubs/ubs.conf.default
@dirrm share/examples/ubs
