# net/autonet/PLIST
@comment $NetBSD: PLIST,v 1.3 2009/06/14 18:09:21 joerg Exp $
man/man8/autonet.8
sbin/autonet
share/examples/autonet/autonet.local
share/examples/autonet/lid_switch
share/examples/rc.d/autonet
