# sysutils/cfengine3/PLIST
@comment $NetBSD: PLIST,v 1.3 2011/02/06 18:44:23 pettai Exp $
lib/libpromises.la
man/man8/cf-agent.8
man/man8/cf-execd.8
man/man8/cf-hub.8
man/man8/cf-key.8
man/man8/cf-know.8
man/man8/cf-monitord.8
man/man8/cf-promises.8
man/man8/cf-report.8
man/man8/cf-runagent.8
man/man8/cf-serverd.8
sbin/cf-agent
sbin/cf-execd
sbin/cf-hub
sbin/cf-key
sbin/cf-know
sbin/cf-monitord
sbin/cf-promises
sbin/cf-report
sbin/cf-runagent
sbin/cf-serverd
share/doc/cfengine/ChangeLog
share/doc/cfengine/INSTALL
share/doc/cfengine/NEWS
share/doc/cfengine/README
share/doc/cfengine/active_directory.cf
share/doc/cfengine/app_baseline.cf
share/doc/cfengine/cfengine_stdlib.cf
share/doc/cfengine/promise_knowledge.cf
share/doc/cfengine/unit_abort.cf
share/doc/cfengine/unit_accessed_before.cf
share/doc/cfengine/unit_accumulated_time.cf
share/doc/cfengine/unit_acl.cf
share/doc/cfengine/unit_acl_generic.cf
share/doc/cfengine/unit_acl_secret.cf
share/doc/cfengine/unit_activedirectory_listusers.cf
share/doc/cfengine/unit_activedirectory_showuser.cf
share/doc/cfengine/unit_ago.cf
share/doc/cfengine/unit_arrays.cf
share/doc/cfengine/unit_backreferences_files.cf
share/doc/cfengine/unit_badpromise.cf
share/doc/cfengine/unit_badtype.cf
share/doc/cfengine/unit_bsdflags.cf
share/doc/cfengine/unit_cf2_integration.cf
share/doc/cfengine/unit_change_detect.cf
share/doc/cfengine/unit_changedbefore.cf
share/doc/cfengine/unit_chdir.cf
share/doc/cfengine/unit_classes_global.cf
share/doc/cfengine/unit_classmatch.cf
share/doc/cfengine/unit_classvar_convergence.cf
share/doc/cfengine/unit_compare.cf
share/doc/cfengine/unit_control_expand.cf
share/doc/cfengine/unit_controlclasses.cf
share/doc/cfengine/unit_copy.cf
share/doc/cfengine/unit_copy_classes.cf
share/doc/cfengine/unit_copy_edit.cf
share/doc/cfengine/unit_copylinks.cf
share/doc/cfengine/unit_countclassesmatching.cf
share/doc/cfengine/unit_countlinesmatching.cf
share/doc/cfengine/unit_create_filedir.cf
share/doc/cfengine/unit_createdb.cf
share/doc/cfengine/unit_customize_by_named_list.cf
share/doc/cfengine/unit_definitions.cf
share/doc/cfengine/unit_deletelines.cf
share/doc/cfengine/unit_disable_and_rotate_files.cf
share/doc/cfengine/unit_diskfree.cf
share/doc/cfengine/unit_dollar.cf
share/doc/cfengine/unit_edit_column_files.cf
share/doc/cfengine/unit_edit_comment_lines.cf
share/doc/cfengine/unit_edit_deletenotmatch.cf
share/doc/cfengine/unit_edit_insert_fuzzylines.cf
share/doc/cfengine/unit_edit_insert_lines.cf
share/doc/cfengine/unit_edit_insert_lines_silly.cf
share/doc/cfengine/unit_edit_passwd_file.cf
share/doc/cfengine/unit_edit_replace_string.cf
share/doc/cfengine/unit_edit_sectioned_file.cf
share/doc/cfengine/unit_edit_setvar.cf
share/doc/cfengine/unit_edit_triggerclass.cf
share/doc/cfengine/unit_env.cf
share/doc/cfengine/unit_environment.cf
share/doc/cfengine/unit_environment_simple.cf
share/doc/cfengine/unit_epimenides.cf
share/doc/cfengine/unit_escape.cf
share/doc/cfengine/unit_exec_args.cf
share/doc/cfengine/unit_exec_in_sequence.cf
share/doc/cfengine/unit_execd.cf
share/doc/cfengine/unit_execresult.cf
share/doc/cfengine/unit_expand.cf
share/doc/cfengine/unit_failedcommand.cf
share/doc/cfengine/unit_failsafe.cf
share/doc/cfengine/unit_file_change_detection.cf
share/doc/cfengine/unit_file_owner_list_template.cf
share/doc/cfengine/unit_fileexists.cf
share/doc/cfengine/unit_filenames.cf
share/doc/cfengine/unit_fileperms.cf
share/doc/cfengine/unit_filesexist.cf
share/doc/cfengine/unit_filesexist2.cf
share/doc/cfengine/unit_filesize.cf
share/doc/cfengine/unit_fix_names.cf
share/doc/cfengine/unit_getenv.cf
share/doc/cfengine/unit_getfields.cf
share/doc/cfengine/unit_getgid.cf
share/doc/cfengine/unit_getindices.cf
share/doc/cfengine/unit_getindixes.cf
share/doc/cfengine/unit_getregistry.cf
share/doc/cfengine/unit_getuid.cf
share/doc/cfengine/unit_getusers.cf
share/doc/cfengine/unit_global_list_expansion.cf
share/doc/cfengine/unit_global_list_expansion_2.cf
share/doc/cfengine/unit_grep.cf
share/doc/cfengine/unit_groupexists.cf
share/doc/cfengine/unit_hash.cf
share/doc/cfengine/unit_hashcomment.cf
share/doc/cfengine/unit_hashmatch.cf
share/doc/cfengine/unit_hashuncomment.cf
share/doc/cfengine/unit_helloworld.cf
share/doc/cfengine/unit_host2ip.cf
share/doc/cfengine/unit_hostrange.cf
share/doc/cfengine/unit_hostsseen.cf
share/doc/cfengine/unit_hub.cf
share/doc/cfengine/unit_id.cf
share/doc/cfengine/unit_insert_users.cf
share/doc/cfengine/unit_intarray.cf
share/doc/cfengine/unit_iprange.cf
share/doc/cfengine/unit_irange.cf
share/doc/cfengine/unit_isdir.cf
share/doc/cfengine/unit_islink.cf
share/doc/cfengine/unit_isnewerthan.cf
share/doc/cfengine/unit_isplain.cf
share/doc/cfengine/unit_isvariable.cf
share/doc/cfengine/unit_iteration.cf
share/doc/cfengine/unit_join.cf
share/doc/cfengine/unit_knowledge_driver.cf
share/doc/cfengine/unit_knowledge_txt.cf
share/doc/cfengine/unit_lastnode.cf
share/doc/cfengine/unit_ldap.cf
share/doc/cfengine/unit_linking.cf
share/doc/cfengine/unit_literal_server.cf
share/doc/cfengine/unit_locate_files_and_compress.cf
share/doc/cfengine/unit_log_private.cf
share/doc/cfengine/unit_loops.cf
share/doc/cfengine/unit_measure_log.cf
share/doc/cfengine/unit_measurements.cf
share/doc/cfengine/unit_menu.cf
share/doc/cfengine/unit_method.cf
share/doc/cfengine/unit_method_validate.cf
share/doc/cfengine/unit_module_exec.cf
share/doc/cfengine/unit_module_exec_2.cf
share/doc/cfengine/unit_monitord.cf
share/doc/cfengine/unit_mount_fs.cf
share/doc/cfengine/unit_multipassvars.cf
share/doc/cfengine/unit_neighbourhood_watch.cf
share/doc/cfengine/unit_null_config.cf
share/doc/cfengine/unit_occurrences.cf
share/doc/cfengine/unit_ordering.cf
share/doc/cfengine/unit_package_apt.cf
share/doc/cfengine/unit_package_freebsd.cf
share/doc/cfengine/unit_package_hash.cf
share/doc/cfengine/unit_package_msi_file.cf
share/doc/cfengine/unit_package_msi_version.cf
share/doc/cfengine/unit_package_rpm.cf
share/doc/cfengine/unit_package_solaris.cf
share/doc/cfengine/unit_package_yum.cf
share/doc/cfengine/unit_package_zypper.cf
share/doc/cfengine/unit_parallel_exec.cf
share/doc/cfengine/unit_pathtype.cf
share/doc/cfengine/unit_pattern_and_edit.cf
share/doc/cfengine/unit_peers.cf
share/doc/cfengine/unit_postfix.cf
share/doc/cfengine/unit_process_kill.cf
share/doc/cfengine/unit_process_matching.cf
share/doc/cfengine/unit_process_matching2.cf
share/doc/cfengine/unit_process_matching3.cf
share/doc/cfengine/unit_process_signalling.cf
share/doc/cfengine/unit_product.cf
share/doc/cfengine/unit_promises.cf
share/doc/cfengine/unit_readlist.cf
share/doc/cfengine/unit_readtcp.cf
share/doc/cfengine/unit_reg_multiline.cf
share/doc/cfengine/unit_regarray.cf
share/doc/cfengine/unit_regex_win.cf
share/doc/cfengine/unit_regextract.cf
share/doc/cfengine/unit_registry.cf
share/doc/cfengine/unit_registry_cache.cf
share/doc/cfengine/unit_regline.cf
share/doc/cfengine/unit_reglist.cf
share/doc/cfengine/unit_remoteclasses.cf
share/doc/cfengine/unit_remove_deadlinks.cf
share/doc/cfengine/unit_rename.cf
share/doc/cfengine/unit_repairedcommand.cf
share/doc/cfengine/unit_report_custom.cf
share/doc/cfengine/unit_report_diffs.cf
share/doc/cfengine/unit_report_state.cf
share/doc/cfengine/unit_reporting.cf
share/doc/cfengine/unit_reporttofile.cf
share/doc/cfengine/unit_resolveconf.cf
share/doc/cfengine/unit_returnszero.cf
share/doc/cfengine/unit_root_passwd.cf
share/doc/cfengine/unit_select_mode.cf
share/doc/cfengine/unit_select_region.cf
share/doc/cfengine/unit_select_size.cf
share/doc/cfengine/unit_selectservers.cf
share/doc/cfengine/unit_server_copy_localhost.cf
share/doc/cfengine/unit_server_copy_purge.cf
share/doc/cfengine/unit_server_copy_remote.cf
share/doc/cfengine/unit_server_flatcopy_localhost.cf
share/doc/cfengine/unit_service_disable.cf
share/doc/cfengine/unit_service_start.cf
share/doc/cfengine/unit_setuidlog.cf
share/doc/cfengine/unit_setvar.cf
share/doc/cfengine/unit_software_dist.cf
share/doc/cfengine/unit_splitstring.cf
share/doc/cfengine/unit_sql.cf
share/doc/cfengine/unit_storage.cf
share/doc/cfengine/unit_strcmp.cf
share/doc/cfengine/unit_stringarray.cf
share/doc/cfengine/unit_switchcase.cf
share/doc/cfengine/unit_syslog.cf
share/doc/cfengine/unit_syslog2.cf
share/doc/cfengine/unit_template.cf
share/doc/cfengine/unit_template2.cf
share/doc/cfengine/unit_test_environment.cf
share/doc/cfengine/unit_tidy_all_files.cf
share/doc/cfengine/unit_translatepath.cf
share/doc/cfengine/unit_update.cf
share/doc/cfengine/unit_user_edit.cf
share/doc/cfengine/unit_user_edit_method.cf
share/doc/cfengine/unit_userexists.cf
share/doc/cfengine/unit_varclass.cf
share/doc/cfengine/unit_varexpansion.cf
share/doc/cfengine/unit_varnet.cf
share/doc/cfengine/unit_vars.cf
share/doc/cfengine/unit_warnifline.cf
share/doc/cfengine/unit_webserver.cf
share/doc/cfengine/unit_zenoss.cf
share/doc/cfengine/win_emergency.cf
share/doc/cfengine/win_registry.cf
share/doc/cfengine/win_schedule.cf
share/doc/cfengine/win_services.cf
share/examples/cfengine/inputs/cfengine_stdlib.cf
share/examples/cfengine/inputs/failsafe.cf
share/examples/cfengine/inputs/promises.cf
share/examples/cfengine/inputs/site.cf
share/examples/cfengine/inputs/update.cf
share/examples/rc.d/cfexecd
share/examples/rc.d/cfmonitord
share/examples/rc.d/cfserverd
@pkgdir etc/cfengine
