# chat/dircproxy/PLIST
@comment $NetBSD: PLIST,v 1.2 2003/06/02 15:20:06 jmmv Exp $
bin/dircproxy
bin/dircproxy-crypt
etc/rc.d/dircproxy
man/man1/dircproxy-crypt.1
man/man1/dircproxy.1
share/dircproxy/cronchk.sh
share/dircproxy/log.pl
share/dircproxy/privmsg-log.pl
share/examples/dircproxy/dircproxyrc
@dirrm share/examples/dircproxy
@dirrm share/dircproxy
