# chat/dircproxy/PLIST
@comment $NetBSD: PLIST,v 1.6 2009/06/14 17:37:43 joerg Exp $
bin/dircproxy
bin/dircproxy-crypt
man/man1/dircproxy-crypt.1
man/man1/dircproxy.1
share/dircproxy/cronchk.sh
share/dircproxy/log.pl
share/dircproxy/privmsg-log.pl
share/examples/dircproxy/dircproxyrc.sample
share/examples/rc.d/dircproxy
