# chat/inspircd12/PLIST
@comment $NetBSD: PLIST,v 1.3 2009/06/06 17:02:21 adrianp Exp $
bin/inspircd
lib/inspircd/cmd_admin.so
lib/inspircd/cmd_away.so
lib/inspircd/cmd_clearcache.so
lib/inspircd/cmd_commands.so
lib/inspircd/cmd_connect.so
lib/inspircd/cmd_die.so
lib/inspircd/cmd_eline.so
lib/inspircd/cmd_gline.so
lib/inspircd/cmd_info.so
lib/inspircd/cmd_invite.so
lib/inspircd/cmd_ison.so
lib/inspircd/cmd_join.so
lib/inspircd/cmd_kick.so
lib/inspircd/cmd_kill.so
lib/inspircd/cmd_kline.so
lib/inspircd/cmd_links.so
lib/inspircd/cmd_list.so
lib/inspircd/cmd_loadmodule.so
lib/inspircd/cmd_lusers.so
lib/inspircd/cmd_map.so
lib/inspircd/cmd_mode.so
lib/inspircd/cmd_modules.so
lib/inspircd/cmd_motd.so
lib/inspircd/cmd_names.so
lib/inspircd/cmd_nick.so
lib/inspircd/cmd_notice.so
lib/inspircd/cmd_oper.so
lib/inspircd/cmd_part.so
lib/inspircd/cmd_pass.so
lib/inspircd/cmd_ping.so
lib/inspircd/cmd_pong.so
lib/inspircd/cmd_privmsg.so
lib/inspircd/cmd_qline.so
lib/inspircd/cmd_quit.so
lib/inspircd/cmd_rehash.so
lib/inspircd/cmd_reloadmodule.so
lib/inspircd/cmd_restart.so
lib/inspircd/cmd_rules.so
lib/inspircd/cmd_server.so
lib/inspircd/cmd_squit.so
lib/inspircd/cmd_stats.so
lib/inspircd/cmd_time.so
lib/inspircd/cmd_topic.so
lib/inspircd/cmd_trace.so
lib/inspircd/cmd_unloadmodule.so
lib/inspircd/cmd_user.so
lib/inspircd/cmd_userhost.so
lib/inspircd/cmd_version.so
lib/inspircd/cmd_wallops.so
lib/inspircd/cmd_who.so
lib/inspircd/cmd_whois.so
lib/inspircd/cmd_whowas.so
lib/inspircd/cmd_zline.so
share/examples/inspircd/aliases/anope.conf.example
share/examples/inspircd/aliases/atheme.conf.example
share/examples/inspircd/aliases/ircservices.conf.example
share/examples/inspircd/aliases/neostats.conf.example
share/examples/inspircd/inspircd.censor.example
share/examples/inspircd/inspircd.conf.example
share/examples/inspircd/inspircd.filter.example
share/examples/inspircd/inspircd.helpop-full.example
share/examples/inspircd/inspircd.helpop.example
share/examples/inspircd/inspircd.motd.example
share/examples/inspircd/inspircd.quotes.example
share/examples/inspircd/inspircd.rules.example
share/examples/inspircd/modules.conf.example
share/examples/inspircd/sql/m_sqllog.mssql.sql
share/examples/inspircd/sql/m_sqllog.mysql.sql
share/examples/inspircd/sql/m_sqllog.postgresql.sql
share/examples/inspircd/sql/m_sqllog.sqlite3.sql
share/examples/inspircd/sql/m_sqloper.mssql.sql
share/examples/inspircd/sql/m_sqloper.mysql.sql
share/examples/inspircd/sql/m_sqloper.postgresql.sql
share/examples/inspircd/sql/m_sqloper.sqlite3.sql
share/examples/rc.d/inspircd
share/inspircd/.gdbargs
share/inspircd/inspircd
share/inspircd/modules/m_alias.so
share/inspircd/modules/m_alltime.so
share/inspircd/modules/m_antibear.so
share/inspircd/modules/m_antibottler.so
share/inspircd/modules/m_auditorium.so
share/inspircd/modules/m_banexception.so
share/inspircd/modules/m_banredirect.so
share/inspircd/modules/m_blockamsg.so
share/inspircd/modules/m_blockcaps.so
share/inspircd/modules/m_blockcolor.so
share/inspircd/modules/m_botmode.so
share/inspircd/modules/m_cban.so
share/inspircd/modules/m_censor.so
share/inspircd/modules/m_cgiirc.so
share/inspircd/modules/m_chancreate.so
share/inspircd/modules/m_chanfilter.so
share/inspircd/modules/m_chanprotect.so
share/inspircd/modules/m_check.so
share/inspircd/modules/m_chghost.so
share/inspircd/modules/m_chgident.so
share/inspircd/modules/m_chgname.so
share/inspircd/modules/m_cloaking.so
share/inspircd/modules/m_clones.so
share/inspircd/modules/m_close.so
share/inspircd/modules/m_commonchans.so
share/inspircd/modules/m_conn_join.so
share/inspircd/modules/m_conn_umodes.so
share/inspircd/modules/m_conn_waitpong.so
share/inspircd/modules/m_connflood.so
share/inspircd/modules/m_customtitle.so
share/inspircd/modules/m_cycle.so
share/inspircd/modules/m_dccallow.so
share/inspircd/modules/m_deaf.so
share/inspircd/modules/m_denychans.so
share/inspircd/modules/m_devoice.so
share/inspircd/modules/m_dnsbl.so
share/inspircd/modules/m_filter.so
share/inspircd/modules/m_foobar.so
share/inspircd/modules/m_globalload.so
share/inspircd/modules/m_globops.so
share/inspircd/modules/m_helpop.so
share/inspircd/modules/m_hidechans.so
share/inspircd/modules/m_hideoper.so
share/inspircd/modules/m_hostchange.so
share/inspircd/modules/m_httpd.so
share/inspircd/modules/m_httpd_stats.so
share/inspircd/modules/m_ident.so
share/inspircd/modules/m_invisible.so
share/inspircd/modules/m_inviteexception.so
share/inspircd/modules/m_joinflood.so
share/inspircd/modules/m_jumpserver.so
share/inspircd/modules/m_kicknorejoin.so
share/inspircd/modules/m_knock.so
share/inspircd/modules/m_lockserv.so
share/inspircd/modules/m_md5.so
share/inspircd/modules/m_messageflood.so
share/inspircd/modules/m_namesx.so
share/inspircd/modules/m_nicklock.so
share/inspircd/modules/m_nickflood.so
share/inspircd/modules/m_noctcp.so
share/inspircd/modules/m_nokicks.so
share/inspircd/modules/m_nonicks.so
share/inspircd/modules/m_nonotice.so
share/inspircd/modules/m_operchans.so
share/inspircd/modules/m_operjoin.so
share/inspircd/modules/m_operlevels.so
share/inspircd/modules/m_operlog.so
share/inspircd/modules/m_opermodes.so
share/inspircd/modules/m_opermotd.so
share/inspircd/modules/m_override.so
share/inspircd/modules/m_randquote.so
share/inspircd/modules/m_redirect.so
share/inspircd/modules/m_regonlycreate.so
share/inspircd/modules/m_remove.so
share/inspircd/modules/m_restrictchans.so
share/inspircd/modules/m_restrictmsg.so
share/inspircd/modules/m_safelist.so
share/inspircd/modules/m_sajoin.so
share/inspircd/modules/m_samode.so
share/inspircd/modules/m_sanick.so
share/inspircd/modules/m_sakick.so
share/inspircd/modules/m_sapart.so
share/inspircd/modules/m_saquit.so
share/inspircd/modules/m_securelist.so
share/inspircd/modules/m_seenicks.so
share/inspircd/modules/m_services_account.so
share/inspircd/modules/m_sethost.so
share/inspircd/modules/m_setident.so
share/inspircd/modules/m_setidle.so
share/inspircd/modules/m_setname.so
share/inspircd/modules/m_sha256.so
share/inspircd/modules/m_showwhois.so
share/inspircd/modules/m_silence.so
share/inspircd/modules/m_spanningtree.so
share/inspircd/modules/m_spy.so
share/inspircd/modules/m_ssl_dummy.so
share/inspircd/modules/m_sslmodes.so
share/inspircd/modules/m_stripcolor.so
share/inspircd/modules/m_svshold.so
share/inspircd/modules/m_swhois.so
share/inspircd/modules/m_taxonomy.so
share/inspircd/modules/m_testcommand.so
share/inspircd/modules/m_timedbans.so
share/inspircd/modules/m_tline.so
share/inspircd/modules/m_uhnames.so
share/inspircd/modules/m_uninvite.so
share/inspircd/modules/m_userip.so
share/inspircd/modules/m_vhost.so
share/inspircd/modules/m_watch.so
share/inspircd/modules/m_abbreviation.so
share/inspircd/modules/m_allowinvite.so
share/inspircd/modules/m_callerid.so
share/inspircd/modules/m_cap.so
share/inspircd/modules/m_chanlog.so
share/inspircd/modules/m_channelban.so
share/inspircd/modules/m_connectban.so
share/inspircd/modules/m_delayjoin.so
share/inspircd/modules/m_gecosban.so
share/inspircd/modules/m_httpd_acl.so
share/inspircd/modules/m_httpd_config.so
share/inspircd/modules/m_maphide.so
share/inspircd/modules/m_muteban.so
share/inspircd/modules/m_nationalchars.so
share/inspircd/modules/m_nopartmsg.so
share/inspircd/modules/m_operinvex.so
share/inspircd/modules/m_operprefix.so
share/inspircd/modules/m_password_hash.so
share/inspircd/modules/m_permchannels.so
share/inspircd/modules/m_regex_glob.so
share/inspircd/modules/m_ripemd160.so
share/inspircd/modules/m_rpc_json.so
share/inspircd/modules/m_rpctest.so
share/inspircd/modules/m_rline.so
share/inspircd/modules/m_sasl.so
share/inspircd/modules/m_satopic.so
share/inspircd/modules/m_serverban.so
share/inspircd/modules/m_servprotect.so
share/inspircd/modules/m_shun.so
share/inspircd/modules/m_xline_db.so
@dirrm lib/inspircd
@dirrm share/examples/inspircd/aliases
@dirrm share/examples/inspircd/sql
@dirrm share/examples/inspircd
