# net/freeradius/PLIST
@comment $NetBSD: PLIST,v 1.24 2009/06/14 18:09:27 joerg Exp $
bin/radclient
bin/radeapclient
bin/radlast
bin/radrelay
bin/radsqlrelay
bin/radtest
bin/radwho
bin/radzap
bin/smbencrypt
lib/libeap-${RADVER}.la
lib/libeap.la
lib/libradius-${RADVER}.la
lib/libradius.la
lib/rlm_acct_unique-${RADVER}.la
lib/rlm_acct_unique.la
lib/rlm_always-${RADVER}.la
lib/rlm_always.la
lib/rlm_attr_filter-${RADVER}.la
lib/rlm_attr_filter.la
lib/rlm_attr_rewrite-${RADVER}.la
lib/rlm_attr_rewrite.la
lib/rlm_chap-${RADVER}.la
lib/rlm_chap.la
lib/rlm_checkval-${RADVER}.la
lib/rlm_checkval.la
lib/rlm_detail-${RADVER}.la
lib/rlm_detail.la
lib/rlm_digest-${RADVER}.la
lib/rlm_digest.la
lib/rlm_eap-${RADVER}.la
lib/rlm_eap.la
lib/rlm_eap_gtc-${RADVER}.la
lib/rlm_eap_gtc.la
lib/rlm_eap_leap-${RADVER}.la
lib/rlm_eap_leap.la
lib/rlm_eap_md5-${RADVER}.la
lib/rlm_eap_md5.la
lib/rlm_eap_mschapv2-${RADVER}.la
lib/rlm_eap_mschapv2.la
lib/rlm_eap_peap-${RADVER}.la
lib/rlm_eap_peap.la
lib/rlm_eap_sim-${RADVER}.la
lib/rlm_eap_sim.la
lib/rlm_eap_tls-${RADVER}.la
lib/rlm_eap_tls.la
lib/rlm_eap_ttls-${RADVER}.la
lib/rlm_eap_ttls.la
lib/rlm_exec-${RADVER}.la
lib/rlm_exec.la
lib/rlm_expr-${RADVER}.la
lib/rlm_expr.la
lib/rlm_fastusers-${RADVER}.la
lib/rlm_fastusers.la
lib/rlm_files-${RADVER}.la
lib/rlm_files.la
lib/rlm_mschap-${RADVER}.la
lib/rlm_mschap.la
lib/rlm_ns_mta_md5-${RADVER}.la
lib/rlm_ns_mta_md5.la
lib/rlm_otp-${RADVER}.la
lib/rlm_otp.la
lib/rlm_pap-${RADVER}.la
lib/rlm_pap.la
lib/rlm_passwd-${RADVER}.la
lib/rlm_passwd.la
lib/rlm_preprocess-${RADVER}.la
lib/rlm_preprocess.la
lib/rlm_radutmp-${RADVER}.la
lib/rlm_radutmp.la
lib/rlm_realm-${RADVER}.la
lib/rlm_realm.la
lib/rlm_sql-${RADVER}.la
lib/rlm_sql.la
lib/rlm_sql_log-${RADVER}.la
lib/rlm_sql_log.la
lib/rlm_sqlcounter-${RADVER}.la
lib/rlm_sqlcounter.la
lib/rlm_sqlippool-${RADVER}.la
lib/rlm_sqlippool.la
lib/rlm_unix-${RADVER}.la
lib/rlm_unix.la
man/man1/radclient.1
man/man1/radeapclient.1
man/man1/radlast.1
man/man1/radtest.1
man/man1/radwho.1
man/man1/radzap.1
man/man5/acct_users.5
man/man5/clients.5
man/man5/clients.conf.5
man/man5/dictionary.5
man/man5/naslist.5
man/man5/radiusd.conf.5
man/man5/rlm_acct_unique.5
man/man5/rlm_always.5
man/man5/rlm_attr_filter.5
man/man5/rlm_attr_rewrite.5
man/man5/rlm_chap.5
man/man5/rlm_counter.5
man/man5/rlm_detail.5
man/man5/rlm_expr.5
man/man5/rlm_files.5
man/man5/rlm_mschap.5
man/man5/rlm_pap.5
man/man5/rlm_passwd.5
man/man5/rlm_realm.5
man/man5/rlm_sql.5
man/man5/rlm_sql_log.5
man/man5/rlm_unix.5
man/man5/users.5
man/man8/radiusd.8
man/man8/radrelay.8
man/man8/radsqlrelay.8
man/man8/radwatch.8
sbin/check-radiusd-config
sbin/checkrad
sbin/radiusd
sbin/radwatch
sbin/rc.radiusd
share/doc/freeradius/Acct-Type
share/doc/freeradius/CYGWIN
share/doc/freeradius/ChangeLog
share/doc/freeradius/Post-Auth-Type
share/doc/freeradius/Session-Type
share/doc/freeradius/ldap_howto.txt
share/doc/freeradius/misc-nas
share/doc/freeradius/release-method.txt
share/doc/freeradius/Autz-Type
share/doc/freeradius/DIFFS
share/doc/freeradius/MACOSX
share/doc/freeradius/OS2
share/doc/freeradius/RADIUS-SQL.schema
share/doc/freeradius/RADIUS-LDAP-eDirectory
share/doc/freeradius/README
share/doc/freeradius/Simultaneous-Use
share/doc/freeradius/aaa.txt
share/doc/freeradius/ascend
share/doc/freeradius/bay
share/doc/freeradius/bugs
share/doc/freeradius/cisco
share/doc/freeradius/coding-methods.txt
share/doc/freeradius/configurable_failover
share/doc/freeradius/duplicate-users
share/doc/freeradius/examples/iplanet.ldif
share/doc/freeradius/examples/iplanet.schema
share/doc/freeradius/examples/mssql.sql
share/doc/freeradius/examples/mysql.sql
share/doc/freeradius/examples/openldap.schema
share/doc/freeradius/examples/oracle.sql
share/doc/freeradius/examples/postgresql.sql
share/doc/freeradius/examples/postgresql.sql.extra
share/doc/freeradius/load-balance.txt
share/doc/freeradius/module_interface
share/doc/freeradius/mssql
share/doc/freeradius/performance-testing
share/doc/freeradius/processing_users_file
share/doc/freeradius/proxy
share/doc/freeradius/radrelay
share/doc/freeradius/rfc/attributes.html
share/doc/freeradius/rfc/draft-kamath-pppext-eap-mschapv2-00.txt
share/doc/freeradius/rfc/draft-sterman-aaa-sip-00.txt
share/doc/freeradius/rfc/leap.txt
share/doc/freeradius/rfc/pppext-eap-sim-12.txt
share/doc/freeradius/rfc/rfc1157.txt
share/doc/freeradius/rfc/rfc1227.txt
share/doc/freeradius/rfc/rfc1448.txt
share/doc/freeradius/rfc/rfc1901.txt
share/doc/freeradius/rfc/rfc1905.txt
share/doc/freeradius/rfc/rfc2058.txt
share/doc/freeradius/rfc/rfc2059.txt
share/doc/freeradius/rfc/rfc2138.txt
share/doc/freeradius/rfc/rfc2139.txt
share/doc/freeradius/rfc/rfc2243.txt
share/doc/freeradius/rfc/rfc2289.txt
share/doc/freeradius/rfc/rfc2433.txt
share/doc/freeradius/rfc/rfc2548.txt
share/doc/freeradius/rfc/rfc2618.txt
share/doc/freeradius/rfc/rfc2619.txt
share/doc/freeradius/rfc/rfc2620.txt
share/doc/freeradius/rfc/rfc2621.txt
share/doc/freeradius/rfc/rfc2716.txt
share/doc/freeradius/rfc/rfc2759.txt
share/doc/freeradius/rfc/rfc2809.txt
share/doc/freeradius/rfc/rfc2865.txt
share/doc/freeradius/rfc/rfc2866.txt
share/doc/freeradius/rfc/rfc2867.txt
share/doc/freeradius/rfc/rfc2868.txt
share/doc/freeradius/rfc/rfc2869.txt
share/doc/freeradius/rfc/rfc2882.txt
share/doc/freeradius/rfc/rfc2924.txt
share/doc/freeradius/rfc/rfc3162.txt
share/doc/freeradius/rfc/rfc3575.txt
share/doc/freeradius/rfc/rfc3576.txt
share/doc/freeradius/rfc/rfc3579.txt
share/doc/freeradius/rfc/rfc3580.txt
share/doc/freeradius/rlm_attr_filter
share/doc/freeradius/rlm_dbm
share/doc/freeradius/rlm_digest
share/doc/freeradius/rlm_eap
share/doc/freeradius/rlm_fastusers
share/doc/freeradius/rlm_krb5
share/doc/freeradius/rlm_ldap
share/doc/freeradius/rlm_pam
share/doc/freeradius/rlm_passwd
share/doc/freeradius/rlm_python
share/doc/freeradius/rlm_sim_triplets
share/doc/freeradius/rlm_sql
share/doc/freeradius/rlm_sqlcounter
share/doc/freeradius/rlm_sqlippool
share/doc/freeradius/supervise-radiusd.txt
share/doc/freeradius/tuning_guide
share/doc/freeradius/variables.txt
share/examples/freeradius/acct_users
share/examples/freeradius/attrs
share/examples/freeradius/certs/README
share/examples/freeradius/certs/cert-clt.der
share/examples/freeradius/certs/cert-clt.p12
share/examples/freeradius/certs/cert-clt.pem
share/examples/freeradius/certs/cert-srv.der
share/examples/freeradius/certs/cert-srv.p12
share/examples/freeradius/certs/cert-srv.pem
share/examples/freeradius/certs/demoCA/cacert.pem
share/examples/freeradius/certs/demoCA/index.txt
share/examples/freeradius/certs/demoCA/index.txt.old
share/examples/freeradius/certs/demoCA/serial
share/examples/freeradius/certs/demoCA/serial.old
share/examples/freeradius/certs/dh
share/examples/freeradius/certs/newcert.pem
share/examples/freeradius/certs/newreq.pem
share/examples/freeradius/certs/random
share/examples/freeradius/certs/root.der
share/examples/freeradius/certs/root.p12
share/examples/freeradius/certs/root.pem
share/examples/freeradius/clients
share/examples/freeradius/clients.conf
share/examples/freeradius/dictionary
share/examples/freeradius/dictionary.in
share/examples/freeradius/eap.conf
share/examples/freeradius/example.pl
share/examples/freeradius/experimental.conf
share/examples/freeradius/hints
share/examples/freeradius/huntgroups
share/examples/freeradius/ldap.attrmap
share/examples/freeradius/mssql.conf
share/examples/freeradius/naslist
share/examples/freeradius/naspasswd
share/examples/freeradius/oraclesql.conf
share/examples/freeradius/otp.conf
share/examples/freeradius/pgsql-voip.conf
share/examples/freeradius/postgresql.conf
share/examples/freeradius/preproxy_users
share/examples/freeradius/proxy.conf
share/examples/freeradius/radiusd.conf
share/examples/freeradius/realms
share/examples/freeradius/snmp.conf
share/examples/freeradius/sql.conf
share/examples/freeradius/sqlippool.conf
share/examples/freeradius/users
share/examples/rc.d/radiusd
share/freeradius/dictionary
share/freeradius/dictionary.3com
share/freeradius/dictionary.3gpp
share/freeradius/dictionary.3gpp2
share/freeradius/dictionary.acc
share/freeradius/dictionary.airespace
share/freeradius/dictionary.alcatel
share/freeradius/dictionary.alteon
share/freeradius/dictionary.altiga
share/freeradius/dictionary.alvarion
share/freeradius/dictionary.aptis
share/freeradius/dictionary.aruba
share/freeradius/dictionary.ascend
share/freeradius/dictionary.asn
share/freeradius/dictionary.avaya
share/freeradius/dictionary.bay
share/freeradius/dictionary.bintec
share/freeradius/dictionary.bristol
share/freeradius/dictionary.cablelabs
share/freeradius/dictionary.cabletron
share/freeradius/dictionary.cisco
share/freeradius/dictionary.cisco.bbsm
share/freeradius/dictionary.cisco.vpn3000
share/freeradius/dictionary.cisco.vpn5000
share/freeradius/dictionary.colubris
share/freeradius/dictionary.columbia_university
share/freeradius/dictionary.compat
share/freeradius/dictionary.cosine
share/freeradius/dictionary.digium
share/freeradius/dictionary.epygi
share/freeradius/dictionary.ericsson
share/freeradius/dictionary.erx
share/freeradius/dictionary.extreme
share/freeradius/dictionary.fortinet
share/freeradius/dictionary.foundry
share/freeradius/dictionary.freeradius
share/freeradius/dictionary.freeradius.internal
share/freeradius/dictionary.gandalf
share/freeradius/dictionary.garderos
share/freeradius/dictionary.gemtek
share/freeradius/dictionary.hp
share/freeradius/dictionary.ipunplugged
share/freeradius/dictionary.issanni
share/freeradius/dictionary.itk
share/freeradius/dictionary.juniper
share/freeradius/dictionary.karlnet
share/freeradius/dictionary.livingston
share/freeradius/dictionary.localweb
share/freeradius/dictionary.lucent
share/freeradius/dictionary.merit
share/freeradius/dictionary.microsoft
share/freeradius/dictionary.mikrotik
share/freeradius/dictionary.motorola
share/freeradius/dictionary.navini
share/freeradius/dictionary.netscreen
share/freeradius/dictionary.nokia
share/freeradius/dictionary.nomadix
share/freeradius/dictionary.nortel
share/freeradius/dictionary.ntua
share/freeradius/dictionary.packeteer
share/freeradius/dictionary.patton
share/freeradius/dictionary.propel
share/freeradius/dictionary.quintum
share/freeradius/dictionary.redback
share/freeradius/dictionary.redcreek
share/freeradius/dictionary.rfc2865
share/freeradius/dictionary.rfc2866
share/freeradius/dictionary.rfc2867
share/freeradius/dictionary.rfc2868
share/freeradius/dictionary.rfc2869
share/freeradius/dictionary.rfc3162
share/freeradius/dictionary.rfc3576
share/freeradius/dictionary.rfc3580
share/freeradius/dictionary.rfc4372
share/freeradius/dictionary.rfc4590
share/freeradius/dictionary.rfc4675
share/freeradius/dictionary.rfc4679
share/freeradius/dictionary.riverstone
share/freeradius/dictionary.roaringpenguin
share/freeradius/dictionary.schulzrinne-sipping
share/freeradius/dictionary.shasta
share/freeradius/dictionary.shiva
share/freeradius/dictionary.sofaware
share/freeradius/dictionary.sonicwall
share/freeradius/dictionary.springtide
share/freeradius/dictionary.starent
share/freeradius/dictionary.t_systems_nova
share/freeradius/dictionary.telebit
share/freeradius/dictionary.trapeze
share/freeradius/dictionary.tropos
share/freeradius/dictionary.unix
share/freeradius/dictionary.usr
share/freeradius/dictionary.valemount
share/freeradius/dictionary.versanet
share/freeradius/dictionary.walabi
share/freeradius/dictionary.waverider
share/freeradius/dictionary.wispr
share/freeradius/dictionary.xedia
share/freeradius/dictionary.xylan
share/freeradius/dictionary.zyxel
@unexec ${RMDIR} %D/etc/raddb/certs/demoCA 2>/dev/null || ${TRUE}
@unexec ${RMDIR} %D/etc/raddb/certs 2>/dev/null || ${TRUE}
