# net/freeradius/PLIST
@comment $NetBSD: PLIST,v 1.8 2005/05/02 20:34:02 reed Exp $
bin/radclient
bin/radeapclient
bin/radlast
bin/radrelay
bin/radtest
bin/radwho
bin/radzap
bin/rlm_dbm_cat
bin/rlm_dbm_parser
bin/smbencrypt
lib/libeap-${PKGVERSION}.la
lib/libeap-${PKGVERSION}.so
lib/libeap.a
lib/libeap.la
lib/libeap.so
lib/libradius-${PKGVERSION}.la
lib/libradius-${PKGVERSION}.so
lib/libradius.a
lib/libradius.la
lib/libradius.so
lib/rlm_acct_unique-${PKGVERSION}.la
lib/rlm_acct_unique-${PKGVERSION}.so
lib/rlm_acct_unique.a
lib/rlm_acct_unique.la
lib/rlm_acct_unique.so
lib/rlm_always-${PKGVERSION}.la
lib/rlm_always-${PKGVERSION}.so
lib/rlm_always.a
lib/rlm_always.la
lib/rlm_always.so
lib/rlm_attr_filter-${PKGVERSION}.la
lib/rlm_attr_filter-${PKGVERSION}.so
lib/rlm_attr_filter.a
lib/rlm_attr_filter.la
lib/rlm_attr_filter.so
lib/rlm_attr_rewrite-${PKGVERSION}.la
lib/rlm_attr_rewrite-${PKGVERSION}.so
lib/rlm_attr_rewrite.a
lib/rlm_attr_rewrite.la
lib/rlm_attr_rewrite.so
lib/rlm_chap-${PKGVERSION}.la
lib/rlm_chap-${PKGVERSION}.so
lib/rlm_chap.a
lib/rlm_chap.la
lib/rlm_chap.so
lib/rlm_checkval-${PKGVERSION}.la
lib/rlm_checkval-${PKGVERSION}.so
lib/rlm_checkval.a
lib/rlm_checkval.la
lib/rlm_checkval.so
lib/rlm_counter-${PKGVERSION}.la
lib/rlm_counter-${PKGVERSION}.so
lib/rlm_counter.a
lib/rlm_counter.la
lib/rlm_counter.so
lib/rlm_dbm-${PKGVERSION}.la
lib/rlm_dbm-${PKGVERSION}.so
lib/rlm_dbm.a
lib/rlm_dbm.la
lib/rlm_dbm.so
lib/rlm_detail-${PKGVERSION}.la
lib/rlm_detail-${PKGVERSION}.so
lib/rlm_detail.a
lib/rlm_detail.la
lib/rlm_detail.so
lib/rlm_digest-${PKGVERSION}.la
lib/rlm_digest-${PKGVERSION}.so
lib/rlm_digest.a
lib/rlm_digest.la
lib/rlm_digest.so
lib/rlm_eap-${PKGVERSION}.la
lib/rlm_eap-${PKGVERSION}.so
lib/rlm_eap.a
lib/rlm_eap.la
lib/rlm_eap.so
lib/rlm_eap_gtc-${PKGVERSION}.la
lib/rlm_eap_gtc-${PKGVERSION}.so
lib/rlm_eap_gtc.a
lib/rlm_eap_gtc.la
lib/rlm_eap_gtc.so
lib/rlm_eap_leap-${PKGVERSION}.la
lib/rlm_eap_leap-${PKGVERSION}.so
lib/rlm_eap_leap.a
lib/rlm_eap_leap.la
lib/rlm_eap_leap.so
lib/rlm_eap_md5-${PKGVERSION}.la
lib/rlm_eap_md5-${PKGVERSION}.so
lib/rlm_eap_md5.a
lib/rlm_eap_md5.la
lib/rlm_eap_md5.so
lib/rlm_eap_mschapv2-${PKGVERSION}.la
lib/rlm_eap_mschapv2-${PKGVERSION}.so
lib/rlm_eap_mschapv2.a
lib/rlm_eap_mschapv2.la
lib/rlm_eap_mschapv2.so
lib/rlm_eap_peap-${PKGVERSION}.la
lib/rlm_eap_peap-${PKGVERSION}.so
lib/rlm_eap_peap.a
lib/rlm_eap_peap.la
lib/rlm_eap_peap.so
lib/rlm_eap_sim-${PKGVERSION}.la
lib/rlm_eap_sim-${PKGVERSION}.so
lib/rlm_eap_sim.a
lib/rlm_eap_sim.la
lib/rlm_eap_sim.so
lib/rlm_eap_tls-${PKGVERSION}.la
lib/rlm_eap_tls-${PKGVERSION}.so
lib/rlm_eap_tls.a
lib/rlm_eap_tls.la
lib/rlm_eap_tls.so
lib/rlm_eap_ttls-${PKGVERSION}.la
lib/rlm_eap_ttls-${PKGVERSION}.so
lib/rlm_eap_ttls.a
lib/rlm_eap_ttls.la
lib/rlm_eap_ttls.so
lib/rlm_exec-${PKGVERSION}.la
lib/rlm_exec-${PKGVERSION}.so
lib/rlm_exec.a
lib/rlm_exec.la
lib/rlm_exec.so
lib/rlm_expr-${PKGVERSION}.la
lib/rlm_expr-${PKGVERSION}.so
lib/rlm_expr.a
lib/rlm_expr.la
lib/rlm_expr.so
lib/rlm_fastusers-${PKGVERSION}.la
lib/rlm_fastusers-${PKGVERSION}.so
lib/rlm_fastusers.a
lib/rlm_fastusers.la
lib/rlm_fastusers.so
lib/rlm_files-${PKGVERSION}.la
lib/rlm_files-${PKGVERSION}.so
lib/rlm_files.a
lib/rlm_files.la
lib/rlm_files.so
lib/rlm_mschap-${PKGVERSION}.la
lib/rlm_mschap-${PKGVERSION}.so
lib/rlm_mschap.a
lib/rlm_mschap.la
lib/rlm_mschap.so
lib/rlm_ns_mta_md5-${PKGVERSION}.la
lib/rlm_ns_mta_md5-${PKGVERSION}.so
lib/rlm_ns_mta_md5.a
lib/rlm_ns_mta_md5.la
lib/rlm_ns_mta_md5.so
lib/rlm_pap-${PKGVERSION}.la
lib/rlm_pap-${PKGVERSION}.so
lib/rlm_pap.a
lib/rlm_pap.la
lib/rlm_pap.so
lib/rlm_passwd-${PKGVERSION}.la
lib/rlm_passwd-${PKGVERSION}.so
lib/rlm_passwd.a
lib/rlm_passwd.la
lib/rlm_passwd.so
lib/rlm_preprocess-${PKGVERSION}.la
lib/rlm_preprocess-${PKGVERSION}.so
lib/rlm_preprocess.a
lib/rlm_preprocess.la
lib/rlm_preprocess.so
lib/rlm_radutmp-${PKGVERSION}.la
lib/rlm_radutmp-${PKGVERSION}.so
lib/rlm_radutmp.a
lib/rlm_radutmp.la
lib/rlm_radutmp.so
lib/rlm_realm-${PKGVERSION}.la
lib/rlm_realm-${PKGVERSION}.so
lib/rlm_realm.a
lib/rlm_realm.la
lib/rlm_realm.so
lib/rlm_sql-${PKGVERSION}.la
lib/rlm_sql-${PKGVERSION}.so
lib/rlm_sql.a
lib/rlm_sql.la
lib/rlm_sql.so
lib/rlm_unix-${PKGVERSION}.la
lib/rlm_unix-${PKGVERSION}.so
lib/rlm_unix.a
lib/rlm_unix.la
lib/rlm_unix.so
lib/rlm_x99_token-${PKGVERSION}.la
lib/rlm_x99_token-${PKGVERSION}.so
lib/rlm_x99_token.a
lib/rlm_x99_token.la
lib/rlm_x99_token.so
man/man1/radclient.1
man/man1/radeapclient.1
man/man1/radlast.1
man/man1/radtest.1
man/man1/radwho.1
man/man1/radzap.1
man/man5/acct_users.5
man/man5/clients.5
man/man5/dictionary.5
man/man5/naslist.5
man/man5/radiusd.conf.5
man/man5/rlm_acct_unique.5
man/man5/rlm_always.5
man/man5/rlm_attr_filter.5
man/man5/rlm_attr_rewrite.5
man/man5/rlm_chap.5
man/man5/rlm_counter.5
man/man5/rlm_detail.5
man/man5/rlm_expr.5
man/man5/rlm_files.5
man/man5/rlm_mschap.5
man/man5/rlm_pap.5
man/man5/rlm_passwd.5
man/man5/rlm_realm.5
man/man5/rlm_sql.5
man/man5/rlm_unix.5
man/man5/users.5
man/man8/radiusd.8
man/man8/radrelay.8
man/man8/radwatch.8
sbin/check-radiusd-config
sbin/checkrad
sbin/radiusd
sbin/radwatch
sbin/rc.radiusd
share/doc/${PKGNAME}/Autz-Type
share/doc/${PKGNAME}/DIFFS
share/doc/${PKGNAME}/MACOSX
share/doc/${PKGNAME}/OS2
share/doc/${PKGNAME}/RADIUS-LDAP.schema
share/doc/${PKGNAME}/RADIUS-LDAPv3.schema
share/doc/${PKGNAME}/RADIUS-SQL.schema
share/doc/${PKGNAME}/README
share/doc/${PKGNAME}/Simultaneous-Use
share/doc/${PKGNAME}/aaa.txt
share/doc/${PKGNAME}/ascend
share/doc/${PKGNAME}/bay
share/doc/${PKGNAME}/bugs
share/doc/${PKGNAME}/cisco
share/doc/${PKGNAME}/coding-methods.txt
share/doc/${PKGNAME}/configurable_failover
share/doc/${PKGNAME}/duplicate-users
share/doc/${PKGNAME}/module_interface
share/doc/${PKGNAME}/performance-testing
share/doc/${PKGNAME}/processing_users_file
share/doc/${PKGNAME}/proxy
share/doc/${PKGNAME}/radrelay
share/doc/${PKGNAME}/rfc/attributes.html
share/doc/${PKGNAME}/rfc/draft-kamath-pppext-eap-mschapv2-00.txt
share/doc/${PKGNAME}/rfc/draft-sterman-aaa-sip-00.txt
share/doc/${PKGNAME}/rfc/leap.txt
share/doc/${PKGNAME}/rfc/pppext-eap-sim-12.txt
share/doc/${PKGNAME}/rfc/rfc1157.txt
share/doc/${PKGNAME}/rfc/rfc1227.txt
share/doc/${PKGNAME}/rfc/rfc1448.txt
share/doc/${PKGNAME}/rfc/rfc1901.txt
share/doc/${PKGNAME}/rfc/rfc1905.txt
share/doc/${PKGNAME}/rfc/rfc2058.txt
share/doc/${PKGNAME}/rfc/rfc2059.txt
share/doc/${PKGNAME}/rfc/rfc2138.txt
share/doc/${PKGNAME}/rfc/rfc2139.txt
share/doc/${PKGNAME}/rfc/rfc2243.txt
share/doc/${PKGNAME}/rfc/rfc2289.txt
share/doc/${PKGNAME}/rfc/rfc2433.txt
share/doc/${PKGNAME}/rfc/rfc2548.txt
share/doc/${PKGNAME}/rfc/rfc2618.txt
share/doc/${PKGNAME}/rfc/rfc2619.txt
share/doc/${PKGNAME}/rfc/rfc2620.txt
share/doc/${PKGNAME}/rfc/rfc2621.txt
share/doc/${PKGNAME}/rfc/rfc2716.txt
share/doc/${PKGNAME}/rfc/rfc2759.txt
share/doc/${PKGNAME}/rfc/rfc2809.txt
share/doc/${PKGNAME}/rfc/rfc2865.txt
share/doc/${PKGNAME}/rfc/rfc2866.txt
share/doc/${PKGNAME}/rfc/rfc2867.txt
share/doc/${PKGNAME}/rfc/rfc2868.txt
share/doc/${PKGNAME}/rfc/rfc2869.txt
share/doc/${PKGNAME}/rfc/rfc2882.txt
share/doc/${PKGNAME}/rfc/rfc2924.txt
share/doc/${PKGNAME}/rfc/rfc3162.txt
share/doc/${PKGNAME}/rfc/rfc3575.txt
share/doc/${PKGNAME}/rfc/rfc3576.txt
share/doc/${PKGNAME}/rfc/rfc3579.txt
share/doc/${PKGNAME}/rfc/rfc3580.txt
share/doc/${PKGNAME}/rlm_attr_filter
share/doc/${PKGNAME}/rlm_dbm
share/doc/${PKGNAME}/rlm_digest
share/doc/${PKGNAME}/rlm_eap
share/doc/${PKGNAME}/rlm_fastusers
share/doc/${PKGNAME}/rlm_krb5
share/doc/${PKGNAME}/rlm_ldap
share/doc/${PKGNAME}/rlm_pam
share/doc/${PKGNAME}/rlm_passwd
share/doc/${PKGNAME}/rlm_python
share/doc/${PKGNAME}/rlm_sim_triplets
share/doc/${PKGNAME}/rlm_sql
share/doc/${PKGNAME}/rlm_sqlcounter
share/doc/${PKGNAME}/rlm_x99_token
share/doc/${PKGNAME}/supervise-radiusd.txt
share/doc/${PKGNAME}/tuning_guide
share/doc/${PKGNAME}/variables.txt
share/examples/freeradius/acct_users
share/examples/freeradius/attrs
share/examples/freeradius/certs/README
share/examples/freeradius/certs/cert-clt.der
share/examples/freeradius/certs/cert-clt.p12
share/examples/freeradius/certs/cert-clt.pem
share/examples/freeradius/certs/cert-srv.der
share/examples/freeradius/certs/cert-srv.p12
share/examples/freeradius/certs/cert-srv.pem
share/examples/freeradius/certs/demoCA/cacert.pem
share/examples/freeradius/certs/demoCA/index.txt
share/examples/freeradius/certs/demoCA/index.txt.old
share/examples/freeradius/certs/demoCA/serial
share/examples/freeradius/certs/demoCA/serial.old
share/examples/freeradius/certs/dh
share/examples/freeradius/certs/newcert.pem
share/examples/freeradius/certs/newreq.pem
share/examples/freeradius/certs/random
share/examples/freeradius/certs/root.der
share/examples/freeradius/certs/root.p12
share/examples/freeradius/certs/root.pem
share/examples/freeradius/clients
share/examples/freeradius/clients.conf
share/examples/freeradius/dictionary
share/examples/freeradius/dictionary.in
share/examples/freeradius/eap.conf
share/examples/freeradius/experimental.conf
share/examples/freeradius/hints
share/examples/freeradius/huntgroups
share/examples/freeradius/ldap.attrmap
share/examples/freeradius/mssql.conf
share/examples/freeradius/naslist
share/examples/freeradius/naspasswd
share/examples/freeradius/oraclesql.conf
share/examples/freeradius/pgsql-voip.conf
share/examples/freeradius/postgresql.conf
share/examples/freeradius/preproxy_users
share/examples/freeradius/proxy.conf
share/examples/freeradius/radiusd.conf
share/examples/freeradius/radiusd.conf.in
share/examples/freeradius/realms
share/examples/freeradius/snmp.conf
share/examples/freeradius/sql.conf
share/examples/freeradius/users
share/examples/freeradius/x99.conf
share/examples/freeradius/x99passwd.sample
share/examples/rc.d/radiusd
share/freeradius/dictionary
share/freeradius/dictionary.3com
share/freeradius/dictionary.3gpp
share/freeradius/dictionary.3gpp2
share/freeradius/dictionary.acc
share/freeradius/dictionary.alcatel
share/freeradius/dictionary.alteon
share/freeradius/dictionary.altiga
share/freeradius/dictionary.aptis
share/freeradius/dictionary.ascend
share/freeradius/dictionary.bay
share/freeradius/dictionary.bintec
share/freeradius/dictionary.bristol
share/freeradius/dictionary.cabletron
share/freeradius/dictionary.cisco
share/freeradius/dictionary.cisco.bbsm
share/freeradius/dictionary.cisco.vpn3000
share/freeradius/dictionary.cisco.vpn5000
share/freeradius/dictionary.colubris
share/freeradius/dictionary.columbia_university
share/freeradius/dictionary.compat
share/freeradius/dictionary.erx
share/freeradius/dictionary.extreme
share/freeradius/dictionary.foundry
share/freeradius/dictionary.freeradius
share/freeradius/dictionary.gandalf
share/freeradius/dictionary.garderos
share/freeradius/dictionary.itk
share/freeradius/dictionary.juniper
share/freeradius/dictionary.karlnet
share/freeradius/dictionary.livingston
share/freeradius/dictionary.localweb
share/freeradius/dictionary.merit
share/freeradius/dictionary.microsoft
share/freeradius/dictionary.mikrotik
share/freeradius/dictionary.navini
share/freeradius/dictionary.netscreen
share/freeradius/dictionary.nokia
share/freeradius/dictionary.nomadix
share/freeradius/dictionary.propel
share/freeradius/dictionary.quintum
share/freeradius/dictionary.redback
share/freeradius/dictionary.redcreek
share/freeradius/dictionary.shasta
share/freeradius/dictionary.shiva
share/freeradius/dictionary.sonicwall
share/freeradius/dictionary.springtide
share/freeradius/dictionary.telebit
share/freeradius/dictionary.trapeze
share/freeradius/dictionary.tunnel
share/freeradius/dictionary.unix
share/freeradius/dictionary.usr
share/freeradius/dictionary.valemount
share/freeradius/dictionary.versanet
share/freeradius/dictionary.wispr
share/freeradius/dictionary.xedia
@dirrm share/freeradius
@dirrm share/examples/freeradius/certs/demoCA
@dirrm share/examples/freeradius/certs
@dirrm share/examples/freeradius
@dirrm share/doc/${PKGNAME}/rfc
@dirrm share/doc/${PKGNAME}
@unexec ${RMDIR} %D/etc/raddb/certs/demoCA 2>/dev/null || ${TRUE}
@unexec ${RMDIR} %D/etc/raddb/certs 2>/dev/null || ${TRUE}
