# mail/exim/PLIST
@comment $NetBSD: PLIST,v 1.6 2003/09/02 10:33:38 abs Exp $
etc/rc.d/exim
man/man8/exim.8
sbin/exicyclog
sbin/exigrep
sbin/exim
sbin/${PKGNAME}-1
sbin/exim_checkaccess
sbin/exim_dbmbuild
sbin/exim_dumpdb
sbin/exim_fixdb
sbin/exim_lock
sbin/exim_tidydb
sbin/eximstats
sbin/exinext
sbin/exiqsumm
sbin/exiwhat
share/examples/exim/aliases
share/examples/exim/configure
share/examples/exim/mailer.conf
@dirrm share/examples/exim
