# mail/exim/PLIST
@comment $NetBSD: PLIST,v 1.11 2005/05/02 20:34:00 reed Exp $
man/man8/exim.8
sbin/${DISTNAME}-1
sbin/exicyclog
sbin/exigrep
sbin/exim
sbin/exim_checkaccess
sbin/exim_dbmbuild
sbin/exim_dumpdb
sbin/exim_fixdb
sbin/exim_lock
sbin/exim_tidydb
sbin/eximstats
sbin/exinext
sbin/exipick
sbin/exiqgrep
sbin/exiqsumm
sbin/exiwhat
share/examples/exim/aliases
share/examples/exim/configure
share/examples/exim/mailer.conf
share/examples/rc.d/exim
@dirrm share/examples/exim
