# mail/exim/PLIST
@comment $NetBSD: PLIST,v 1.4 2002/12/09 11:38:04 ad Exp $
etc/exim/configure.example
etc/exim/mailer.conf.exim
etc/rc.d/exim
sbin/exicyclog
sbin/exigrep
sbin/exim
sbin/${PKGNAME}-1
sbin/exim_checkaccess
sbin/exim_dbmbuild
sbin/exim_dumpdb
sbin/exim_fixdb
sbin/exim_lock
sbin/exim_newaliases
sbin/exim_tidydb
sbin/eximstats
sbin/exinext
sbin/exiqsumm
sbin/exiwhat
@dirrm etc/exim
