# mail/exim/PLIST
@comment $NetBSD: PLIST,v 1.5 2003/08/02 17:00:08 jmmv Exp $
etc/rc.d/exim
sbin/exicyclog
sbin/exigrep
sbin/exim
sbin/${PKGNAME}-1
sbin/exim_checkaccess
sbin/exim_dbmbuild
sbin/exim_dumpdb
sbin/exim_fixdb
sbin/exim_lock
sbin/exim_newaliases
sbin/exim_tidydb
sbin/eximstats
sbin/exinext
sbin/exiqsumm
sbin/exiwhat
share/examples/exim/aliases
share/examples/exim/configure
share/examples/exim/mailer.conf
@dirrm share/examples/exim
