# mail/exim/PLIST
@comment $NetBSD: PLIST,v 1.3 2002/06/19 16:02:26 ad Exp $
etc/exim/configure.example
etc/exim/mailer.conf.exim
etc/rc.d/exim
sbin/exicyclog
sbin/exigrep
sbin/exim
sbin/exim-4.05-1
sbin/exim_checkaccess
sbin/exim_dbmbuild
sbin/exim_dumpdb
sbin/exim_fixdb
sbin/exim_lock
sbin/exim_newaliases
sbin/exim_tidydb
sbin/eximstats
sbin/exinext
sbin/exiqsumm
sbin/exiwhat
@unexec ${RMDIR} /var/log/exim 2>/dev/null || true
@dirrm etc/exim
