# mail/nullmailer/PLIST
@comment $NetBSD: PLIST,v 1.7 2005/05/02 20:34:00 reed Exp $
libexec/nullmailer/mailq
share/examples/rc.d/nullmailer
@exec ${CHMOD} u+s %D/libexec/nullmailer/mailq
@exec ${CHOWN} ${NULLMAILER_USER}:${NULLMAILER_GROUP} %D/libexec/nullmailer/mailq
libexec/nullmailer/nullmailer-inject
libexec/nullmailer/nullmailer-queue
@exec ${CHMOD} u+s %D/libexec/nullmailer/nullmailer-queue
@exec ${CHOWN} ${NULLMAILER_USER}:${NULLMAILER_GROUP} %D/libexec/nullmailer/nullmailer-queue
libexec/nullmailer/nullmailer-send
libexec/nullmailer/qmqp
libexec/nullmailer/sendmail
libexec/nullmailer/smtp
man/man1/nullmailer-inject.1
man/man1/sendmail.1
man/man7/nullmailer.7
man/man8/nullmailer-queue.8
man/man8/nullmailer-send.8
share/doc/nullmailer/AUTHORS
share/doc/nullmailer/BUGS
share/doc/nullmailer/COPYING
share/doc/nullmailer/ChangeLog
share/doc/nullmailer/HOWTO
share/doc/nullmailer/NEWS
share/doc/nullmailer/README
share/doc/nullmailer/TODO
share/examples/nullmailer/mailer.conf
@exec ${MKDIR} /var/spool/nullmailer
@exec ${MKDIR} /var/spool/nullmailer/queue
@exec ${MKDIR} /var/spool/nullmailer/tmp
@exec mkfifo /var/spool/nullmailer/trigger
@exec ${CHMOD} 600 /var/spool/nullmailer/trigger
@exec ${CHOWN} -R ${NULLMAILER_USER}:${NULLMAILER_GROUP} /var/spool/nullmailer/*
@dirrm libexec/nullmailer
@dirrm share/examples/nullmailer
@dirrm share/doc/nullmailer
@unexec ${RMDIR} ${PKG_SYSCONFDIR}/nullmailer 2>/dev/null || ${TRUE}
@unexec ${RMDIR} /var/spool/nullmailer/queue 2>/dev/null || ${TRUE}
@unexec ${RMDIR} /var/spool/nullmailer/tmp 2>/dev/null || ${TRUE}
@unexec ${RM} -f /var/spool/nullmailer/trigger 2>/dev/null || ${TRUE}
@unexec ${RMDIR} /var/spool/nullmailer 2>/dev/null || ${TRUE}
