# sysutils/estd/PLIST
@comment $NetBSD: PLIST,v 1.2 2005/12/27 09:58:56 dillo Exp $
sbin/estd
share/examples/rc.d/estd
man/man1/estd.1
