# www/squid26/PLIST
@comment $NetBSD: PLIST,v 1.1.1.1 2008/09/12 15:09:19 taca Exp $
${PLIST.coss}bin/cossdump
bin/squidclient
libexec/cachemgr.cgi
${PLIST.da_ldap}libexec/digest_ldap_auth
${PLIST.da_password}libexec/digest_pw_auth
${PLIST.diskd}libexec/diskd-daemon
${PLIST.na_fakeauth}libexec/fakeauth_auth
${PLIST.ba_getpwnam}libexec/getpwname_auth
${PLIST.eacl_ip_user}libexec/ip_user_check
${PLIST.ba_MSNT}libexec/msnt_auth
${PLIST.ba_NCSA}libexec/ncsa_auth
${PLIST.na_SMB}libexec/ntlm_auth
${PLIST.ba_PAM}libexec/pam_auth
libexec/pinger
${PLIST.ba_LDAP}libexec/squid_ldap_auth
${PLIST.eacl_ldap_group}libexec/squid_ldap_group
${PLIST.eacl_unix_group}libexec/squid_unix_group
${PLIST.unlinkd}libexec/unlinkd
${PLIST.ba_YP}libexec/yp_auth
man/man8/cachemgr.cgi.8
${PLIST.ba_NCSA}man/man8/ncsa_auth.8
${PLIST.ba_PAM}man/man8/pam_auth.8
man/man8/squid.8
${PLIST.ba_LDAP}man/man8/squid_ldap_auth.8
${PLIST.eacl_ldap_group}man/man8/squid_ldap_group.8
${PLIST.eacl_unix_group}man/man8/squid_unix_group.8
sbin/squid
share/doc/squid/ChangeLog
share/doc/squid/README.FreeBSD
share/doc/squid/README.NetBSD
share/doc/squid/README.OpenBSD
share/doc/squid/README.Solaris
share/doc/squid/RELEASENOTES.html
share/doc/squid/debug-sections.txt
share/examples/rc.d/squid
share/examples/squid/cachemgr.conf
share/examples/squid/mime.conf
share/examples/squid/msntauth.conf
share/examples/squid/squid.conf
${PLIST.snmp}share/squid/mib.txt
