# chat/ejabberd/PLIST
@comment $NetBSD: PLIST,v 1.8 2009/06/14 18:28:43 joerg Exp $
${EJABBERD_EXDIR}/ejabberd.cfg
${EJABBERD_EXDIR}/ejabberdctl.cfg
${EJABBERD_EXDIR}/inetrc
lib/erlang/lib/${PKGNAME}/ebin/ELDAPv3.beam
lib/erlang/lib/${PKGNAME}/ebin/XmppAddr.beam
lib/erlang/lib/${PKGNAME}/ebin/acl.beam
lib/erlang/lib/${PKGNAME}/ebin/adhoc.beam
lib/erlang/lib/${PKGNAME}/ebin/cyrsasl.beam
lib/erlang/lib/${PKGNAME}/ebin/cyrsasl_anonymous.beam
lib/erlang/lib/${PKGNAME}/ebin/cyrsasl_digest.beam
lib/erlang/lib/${PKGNAME}/ebin/cyrsasl_plain.beam
lib/erlang/lib/${PKGNAME}/ebin/ejabberd.app
lib/erlang/lib/${PKGNAME}/ebin/ejabberd.beam
lib/erlang/lib/${PKGNAME}/ebin/ejabberd_admin.beam
lib/erlang/lib/${PKGNAME}/ebin/ejabberd_app.beam
lib/erlang/lib/${PKGNAME}/ebin/ejabberd_auth.beam
lib/erlang/lib/${PKGNAME}/ebin/ejabberd_auth_anonymous.beam
lib/erlang/lib/${PKGNAME}/ebin/ejabberd_auth_external.beam
lib/erlang/lib/${PKGNAME}/ebin/ejabberd_auth_internal.beam
lib/erlang/lib/${PKGNAME}/ebin/ejabberd_auth_ldap.beam
lib/erlang/lib/${PKGNAME}/ebin/ejabberd_auth_odbc.beam
lib/erlang/lib/${PKGNAME}/ebin/ejabberd_auth_pam.beam
lib/erlang/lib/${PKGNAME}/ebin/ejabberd_c2s.beam
lib/erlang/lib/${PKGNAME}/ebin/ejabberd_c2s_config.beam
lib/erlang/lib/${PKGNAME}/ebin/ejabberd_check.beam
lib/erlang/lib/${PKGNAME}/ebin/ejabberd_config.beam
lib/erlang/lib/${PKGNAME}/ebin/ejabberd_ctl.beam
lib/erlang/lib/${PKGNAME}/ebin/ejabberd_frontend_socket.beam
lib/erlang/lib/${PKGNAME}/ebin/ejabberd_hooks.beam
lib/erlang/lib/${PKGNAME}/ebin/ejabberd_http.beam
lib/erlang/lib/${PKGNAME}/ebin/ejabberd_http_bind.beam
lib/erlang/lib/${PKGNAME}/ebin/ejabberd_http_poll.beam
lib/erlang/lib/${PKGNAME}/ebin/ejabberd_listener.beam
lib/erlang/lib/${PKGNAME}/ebin/ejabberd_local.beam
lib/erlang/lib/${PKGNAME}/ebin/ejabberd_logger_h.beam
lib/erlang/lib/${PKGNAME}/ebin/ejabberd_loglevel.beam
lib/erlang/lib/${PKGNAME}/ebin/ejabberd_node_groups.beam
lib/erlang/lib/${PKGNAME}/ebin/ejabberd_odbc.beam
lib/erlang/lib/${PKGNAME}/ebin/ejabberd_odbc_sup.beam
lib/erlang/lib/${PKGNAME}/ebin/ejabberd_rdbms.beam
lib/erlang/lib/${PKGNAME}/ebin/ejabberd_receiver.beam
lib/erlang/lib/${PKGNAME}/ebin/ejabberd_router.beam
lib/erlang/lib/${PKGNAME}/ebin/ejabberd_s2s.beam
lib/erlang/lib/${PKGNAME}/ebin/ejabberd_s2s_in.beam
lib/erlang/lib/${PKGNAME}/ebin/ejabberd_s2s_out.beam
lib/erlang/lib/${PKGNAME}/ebin/ejabberd_service.beam
lib/erlang/lib/${PKGNAME}/ebin/ejabberd_sm.beam
lib/erlang/lib/${PKGNAME}/ebin/ejabberd_socket.beam
lib/erlang/lib/${PKGNAME}/ebin/ejabberd_sup.beam
lib/erlang/lib/${PKGNAME}/ebin/ejabberd_system_monitor.beam
lib/erlang/lib/${PKGNAME}/ebin/ejabberd_tmp_sup.beam
lib/erlang/lib/${PKGNAME}/ebin/ejabberd_update.beam
lib/erlang/lib/${PKGNAME}/ebin/ejabberd_web.beam
lib/erlang/lib/${PKGNAME}/ebin/ejabberd_web_admin.beam
lib/erlang/lib/${PKGNAME}/ebin/ejabberd_zlib.beam
lib/erlang/lib/${PKGNAME}/ebin/ejd2odbc.beam
lib/erlang/lib/${PKGNAME}/ebin/eldap.beam
lib/erlang/lib/${PKGNAME}/ebin/eldap_filter.beam
lib/erlang/lib/${PKGNAME}/ebin/eldap_pool.beam
lib/erlang/lib/${PKGNAME}/ebin/eldap_utils.beam
lib/erlang/lib/${PKGNAME}/ebin/extauth.beam
lib/erlang/lib/${PKGNAME}/ebin/gen_iq_handler.beam
lib/erlang/lib/${PKGNAME}/ebin/gen_mod.beam
lib/erlang/lib/${PKGNAME}/ebin/gen_pubsub_node.beam
lib/erlang/lib/${PKGNAME}/ebin/gen_pubsub_nodetree.beam
lib/erlang/lib/${PKGNAME}/ebin/iconv.beam
lib/erlang/lib/${PKGNAME}/ebin/idna.beam
lib/erlang/lib/${PKGNAME}/ebin/jd2ejd.beam
lib/erlang/lib/${PKGNAME}/ebin/jlib.beam
lib/erlang/lib/${PKGNAME}/ebin/mod_adhoc.beam
lib/erlang/lib/${PKGNAME}/ebin/mod_announce.beam
lib/erlang/lib/${PKGNAME}/ebin/mod_caps.beam
lib/erlang/lib/${PKGNAME}/ebin/mod_configure.beam
lib/erlang/lib/${PKGNAME}/ebin/mod_configure2.beam
lib/erlang/lib/${PKGNAME}/ebin/mod_disco.beam
lib/erlang/lib/${PKGNAME}/ebin/mod_echo.beam
lib/erlang/lib/${PKGNAME}/ebin/mod_http_bind.beam
lib/erlang/lib/${PKGNAME}/ebin/mod_http_fileserver.beam
lib/erlang/lib/${PKGNAME}/ebin/mod_ip_blacklist.beam
lib/erlang/lib/${PKGNAME}/ebin/mod_irc.beam
lib/erlang/lib/${PKGNAME}/ebin/mod_irc_connection.beam
lib/erlang/lib/${PKGNAME}/ebin/mod_last.beam
lib/erlang/lib/${PKGNAME}/ebin/mod_last_odbc.beam
lib/erlang/lib/${PKGNAME}/ebin/mod_muc.beam
lib/erlang/lib/${PKGNAME}/ebin/mod_muc_log.beam
lib/erlang/lib/${PKGNAME}/ebin/mod_muc_room.beam
lib/erlang/lib/${PKGNAME}/ebin/mod_offline.beam
lib/erlang/lib/${PKGNAME}/ebin/mod_offline_odbc.beam
lib/erlang/lib/${PKGNAME}/ebin/mod_privacy.beam
lib/erlang/lib/${PKGNAME}/ebin/mod_privacy_odbc.beam
lib/erlang/lib/${PKGNAME}/ebin/mod_private.beam
lib/erlang/lib/${PKGNAME}/ebin/mod_private_odbc.beam
lib/erlang/lib/${PKGNAME}/ebin/mod_proxy65.beam
lib/erlang/lib/${PKGNAME}/ebin/mod_proxy65_lib.beam
lib/erlang/lib/${PKGNAME}/ebin/mod_proxy65_service.beam
lib/erlang/lib/${PKGNAME}/ebin/mod_proxy65_sm.beam
lib/erlang/lib/${PKGNAME}/ebin/mod_proxy65_stream.beam
lib/erlang/lib/${PKGNAME}/ebin/mod_pubsub.beam
lib/erlang/lib/${PKGNAME}/ebin/mod_register.beam
lib/erlang/lib/${PKGNAME}/ebin/mod_roster.beam
lib/erlang/lib/${PKGNAME}/ebin/mod_roster_odbc.beam
lib/erlang/lib/${PKGNAME}/ebin/mod_service_log.beam
lib/erlang/lib/${PKGNAME}/ebin/mod_shared_roster.beam
lib/erlang/lib/${PKGNAME}/ebin/mod_stats.beam
lib/erlang/lib/${PKGNAME}/ebin/mod_time.beam
lib/erlang/lib/${PKGNAME}/ebin/mod_vcard.beam
lib/erlang/lib/${PKGNAME}/ebin/mod_vcard_ldap.beam
lib/erlang/lib/${PKGNAME}/ebin/mod_vcard_odbc.beam
lib/erlang/lib/${PKGNAME}/ebin/mod_version.beam
lib/erlang/lib/${PKGNAME}/ebin/node_buddy.beam
lib/erlang/lib/${PKGNAME}/ebin/node_club.beam
lib/erlang/lib/${PKGNAME}/ebin/node_default.beam
lib/erlang/lib/${PKGNAME}/ebin/node_dispatch.beam
lib/erlang/lib/${PKGNAME}/ebin/node_pep.beam
lib/erlang/lib/${PKGNAME}/ebin/node_private.beam
lib/erlang/lib/${PKGNAME}/ebin/node_public.beam
lib/erlang/lib/${PKGNAME}/ebin/nodetree_default.beam
lib/erlang/lib/${PKGNAME}/ebin/nodetree_virtual.beam
lib/erlang/lib/${PKGNAME}/ebin/odbc_queries.beam
lib/erlang/lib/${PKGNAME}/ebin/p1_fsm.beam
lib/erlang/lib/${PKGNAME}/ebin/p1_mnesia.beam
lib/erlang/lib/${PKGNAME}/ebin/ram_file_io_server.beam
lib/erlang/lib/${PKGNAME}/ebin/randoms.beam
lib/erlang/lib/${PKGNAME}/ebin/sha.beam
lib/erlang/lib/${PKGNAME}/ebin/shaper.beam
lib/erlang/lib/${PKGNAME}/ebin/stringprep.beam
lib/erlang/lib/${PKGNAME}/ebin/stringprep_sup.beam
lib/erlang/lib/${PKGNAME}/ebin/tls.beam
lib/erlang/lib/${PKGNAME}/ebin/translate.beam
lib/erlang/lib/${PKGNAME}/ebin/treap.beam
lib/erlang/lib/${PKGNAME}/ebin/win32_dns.beam
lib/erlang/lib/${PKGNAME}/ebin/xml.beam
lib/erlang/lib/${PKGNAME}/ebin/xml_stream.beam
lib/erlang/lib/${PKGNAME}/priv/lib/ejabberd_zlib_drv.so
lib/erlang/lib/${PKGNAME}/priv/lib/expat_erl.so
lib/erlang/lib/${PKGNAME}/priv/lib/iconv_erl.so
lib/erlang/lib/${PKGNAME}/priv/lib/stringprep_drv.so
lib/erlang/lib/${PKGNAME}/priv/lib/tls_drv.so
lib/erlang/lib/${PKGNAME}/priv/msgs/ca.msg
lib/erlang/lib/${PKGNAME}/priv/msgs/cs.msg
lib/erlang/lib/${PKGNAME}/priv/msgs/de.msg
lib/erlang/lib/${PKGNAME}/priv/msgs/eo.msg
lib/erlang/lib/${PKGNAME}/priv/msgs/es.msg
lib/erlang/lib/${PKGNAME}/priv/msgs/fr.msg
lib/erlang/lib/${PKGNAME}/priv/msgs/gl.msg
lib/erlang/lib/${PKGNAME}/priv/msgs/it.msg
lib/erlang/lib/${PKGNAME}/priv/msgs/ja.msg
lib/erlang/lib/${PKGNAME}/priv/msgs/nl.msg
lib/erlang/lib/${PKGNAME}/priv/msgs/no.msg
lib/erlang/lib/${PKGNAME}/priv/msgs/pl.msg
lib/erlang/lib/${PKGNAME}/priv/msgs/pt-br.msg
lib/erlang/lib/${PKGNAME}/priv/msgs/pt.msg
lib/erlang/lib/${PKGNAME}/priv/msgs/ru.msg
lib/erlang/lib/${PKGNAME}/priv/msgs/sk.msg
lib/erlang/lib/${PKGNAME}/priv/msgs/sv.msg
lib/erlang/lib/${PKGNAME}/priv/msgs/th.msg
lib/erlang/lib/${PKGNAME}/priv/msgs/tr.msg
lib/erlang/lib/${PKGNAME}/priv/msgs/uk.msg
lib/erlang/lib/${PKGNAME}/priv/msgs/vi.msg
lib/erlang/lib/${PKGNAME}/priv/msgs/wa.msg
lib/erlang/lib/${PKGNAME}/priv/msgs/zh.msg
sbin/ejabberdctl
share/doc/${PKGNAME}/COPYING
share/doc/${PKGNAME}/Makefile
share/doc/${PKGNAME}/api/Makefile
share/doc/${PKGNAME}/api/overview.edoc
share/doc/${PKGNAME}/api/process-one.css
share/doc/${PKGNAME}/dev.html
share/doc/${PKGNAME}/dev.tex
share/doc/${PKGNAME}/discorus.png
share/doc/${PKGNAME}/features.html
share/doc/${PKGNAME}/features.pdf
share/doc/${PKGNAME}/features.tex
share/doc/${PKGNAME}/flow.dot
share/doc/${PKGNAME}/guide.html
share/doc/${PKGNAME}/guide.pdf
share/doc/${PKGNAME}/guide.tex
share/doc/${PKGNAME}/introduction.tex
share/doc/${PKGNAME}/logo.png
share/doc/${PKGNAME}/mod_http_bind.tex
share/doc/${PKGNAME}/mod_http_fileserver.tex
share/doc/${PKGNAME}/release_notes_${PKGVERSION}.txt
share/doc/${PKGNAME}/release_notes_0.9.1.txt
share/doc/${PKGNAME}/release_notes_0.9.8.txt
share/doc/${PKGNAME}/release_notes_0.9.txt
share/doc/${PKGNAME}/release_notes_1.0.0.txt
share/doc/${PKGNAME}/release_notes_1.1.0.txt
share/doc/${PKGNAME}/release_notes_1.1.1.txt
share/doc/${PKGNAME}/release_notes_1.1.2.txt
share/doc/${PKGNAME}/release_notes_1.1.3.txt
share/doc/${PKGNAME}/release_notes_1.1.4.txt
share/doc/${PKGNAME}/release_notes_2.0.0.txt
share/doc/${PKGNAME}/release_notes_2.0.1.txt
share/doc/${PKGNAME}/release_notes_2.0.2.txt
share/doc/${PKGNAME}/release_notes_2.0.3.txt
share/doc/${PKGNAME}/release_notes_2.0.4.txt
share/doc/${PKGNAME}/version.tex
share/doc/${PKGNAME}/webadmmain.png
share/doc/${PKGNAME}/webadmmainru.png
share/doc/${PKGNAME}/yozhikheader.png
share/examples/ejabberd/ejabberd.cfg
share/examples/ejabberd/ejabberdctl.cfg
share/examples/ejabberd/inetrc
share/examples/rc.d/ejabberd
