# chat/ejabberd/PLIST
@comment $NetBSD: PLIST,v 1.1.1.1.2.1 2008/04/25 01:09:18 rtr Exp $
${EJABBERD_EXDIR}/ejabberd.cfg
${EJABBERD_EXDIR}/ejabberdctl.cfg
${EJABBERD_EXDIR}/inetrc
lib/erlang/lib/${DISTNAME}/ebin/ELDAPv3.beam
lib/erlang/lib/${DISTNAME}/ebin/XmppAddr.beam
lib/erlang/lib/${DISTNAME}/ebin/acl.beam
lib/erlang/lib/${DISTNAME}/ebin/adhoc.beam
lib/erlang/lib/${DISTNAME}/ebin/cyrsasl.beam
lib/erlang/lib/${DISTNAME}/ebin/cyrsasl_anonymous.beam
lib/erlang/lib/${DISTNAME}/ebin/cyrsasl_digest.beam
lib/erlang/lib/${DISTNAME}/ebin/cyrsasl_plain.beam
lib/erlang/lib/${DISTNAME}/ebin/ejabberd.app
lib/erlang/lib/${DISTNAME}/ebin/ejabberd.beam
lib/erlang/lib/${DISTNAME}/ebin/ejabberd_admin.beam
lib/erlang/lib/${DISTNAME}/ebin/ejabberd_app.beam
lib/erlang/lib/${DISTNAME}/ebin/ejabberd_auth.beam
lib/erlang/lib/${DISTNAME}/ebin/ejabberd_auth_anonymous.beam
lib/erlang/lib/${DISTNAME}/ebin/ejabberd_auth_external.beam
lib/erlang/lib/${DISTNAME}/ebin/ejabberd_auth_internal.beam
lib/erlang/lib/${DISTNAME}/ebin/ejabberd_auth_ldap.beam
lib/erlang/lib/${DISTNAME}/ebin/ejabberd_auth_odbc.beam
lib/erlang/lib/${DISTNAME}/ebin/ejabberd_auth_pam.beam
lib/erlang/lib/${DISTNAME}/ebin/ejabberd_c2s.beam
lib/erlang/lib/${DISTNAME}/ebin/ejabberd_c2s_config.beam
lib/erlang/lib/${DISTNAME}/ebin/ejabberd_config.beam
lib/erlang/lib/${DISTNAME}/ebin/ejabberd_ctl.beam
lib/erlang/lib/${DISTNAME}/ebin/ejabberd_frontend_socket.beam
lib/erlang/lib/${DISTNAME}/ebin/ejabberd_hooks.beam
lib/erlang/lib/${DISTNAME}/ebin/ejabberd_http.beam
lib/erlang/lib/${DISTNAME}/ebin/ejabberd_http_bind.beam
lib/erlang/lib/${DISTNAME}/ebin/ejabberd_http_poll.beam
lib/erlang/lib/${DISTNAME}/ebin/ejabberd_listener.beam
lib/erlang/lib/${DISTNAME}/ebin/ejabberd_local.beam
lib/erlang/lib/${DISTNAME}/ebin/ejabberd_logger_h.beam
lib/erlang/lib/${DISTNAME}/ebin/ejabberd_loglevel.beam
lib/erlang/lib/${DISTNAME}/ebin/ejabberd_node_groups.beam
lib/erlang/lib/${DISTNAME}/ebin/ejabberd_odbc.beam
lib/erlang/lib/${DISTNAME}/ebin/ejabberd_odbc_sup.beam
lib/erlang/lib/${DISTNAME}/ebin/ejabberd_rdbms.beam
lib/erlang/lib/${DISTNAME}/ebin/ejabberd_receiver.beam
lib/erlang/lib/${DISTNAME}/ebin/ejabberd_router.beam
lib/erlang/lib/${DISTNAME}/ebin/ejabberd_s2s.beam
lib/erlang/lib/${DISTNAME}/ebin/ejabberd_s2s_in.beam
lib/erlang/lib/${DISTNAME}/ebin/ejabberd_s2s_out.beam
lib/erlang/lib/${DISTNAME}/ebin/ejabberd_service.beam
lib/erlang/lib/${DISTNAME}/ebin/ejabberd_sm.beam
lib/erlang/lib/${DISTNAME}/ebin/ejabberd_socket.beam
lib/erlang/lib/${DISTNAME}/ebin/ejabberd_sup.beam
lib/erlang/lib/${DISTNAME}/ebin/ejabberd_system_monitor.beam
lib/erlang/lib/${DISTNAME}/ebin/ejabberd_tmp_sup.beam
lib/erlang/lib/${DISTNAME}/ebin/ejabberd_update.beam
lib/erlang/lib/${DISTNAME}/ebin/ejabberd_web.beam
lib/erlang/lib/${DISTNAME}/ebin/ejabberd_web_admin.beam
lib/erlang/lib/${DISTNAME}/ebin/ejabberd_zlib.beam
lib/erlang/lib/${DISTNAME}/ebin/ejd2odbc.beam
lib/erlang/lib/${DISTNAME}/ebin/eldap.beam
lib/erlang/lib/${DISTNAME}/ebin/eldap_filter.beam
lib/erlang/lib/${DISTNAME}/ebin/eldap_pool.beam
lib/erlang/lib/${DISTNAME}/ebin/eldap_utils.beam
lib/erlang/lib/${DISTNAME}/ebin/extauth.beam
lib/erlang/lib/${DISTNAME}/ebin/gen_iq_handler.beam
lib/erlang/lib/${DISTNAME}/ebin/gen_mod.beam
lib/erlang/lib/${DISTNAME}/ebin/gen_pubsub_node.beam
lib/erlang/lib/${DISTNAME}/ebin/gen_pubsub_nodetree.beam
lib/erlang/lib/${DISTNAME}/ebin/iconv.beam
lib/erlang/lib/${DISTNAME}/ebin/idna.beam
lib/erlang/lib/${DISTNAME}/ebin/jd2ejd.beam
lib/erlang/lib/${DISTNAME}/ebin/jlib.beam
lib/erlang/lib/${DISTNAME}/ebin/mod_adhoc.beam
lib/erlang/lib/${DISTNAME}/ebin/mod_announce.beam
lib/erlang/lib/${DISTNAME}/ebin/mod_caps.beam
lib/erlang/lib/${DISTNAME}/ebin/mod_configure.beam
lib/erlang/lib/${DISTNAME}/ebin/mod_configure2.beam
lib/erlang/lib/${DISTNAME}/ebin/mod_disco.beam
lib/erlang/lib/${DISTNAME}/ebin/mod_echo.beam
lib/erlang/lib/${DISTNAME}/ebin/mod_http_bind.beam
lib/erlang/lib/${DISTNAME}/ebin/mod_http_fileserver.beam
lib/erlang/lib/${DISTNAME}/ebin/mod_irc.beam
lib/erlang/lib/${DISTNAME}/ebin/mod_irc_connection.beam
lib/erlang/lib/${DISTNAME}/ebin/mod_last.beam
lib/erlang/lib/${DISTNAME}/ebin/mod_last_odbc.beam
lib/erlang/lib/${DISTNAME}/ebin/mod_muc.beam
lib/erlang/lib/${DISTNAME}/ebin/mod_muc_log.beam
lib/erlang/lib/${DISTNAME}/ebin/mod_muc_room.beam
lib/erlang/lib/${DISTNAME}/ebin/mod_offline.beam
lib/erlang/lib/${DISTNAME}/ebin/mod_offline_odbc.beam
lib/erlang/lib/${DISTNAME}/ebin/mod_privacy.beam
lib/erlang/lib/${DISTNAME}/ebin/mod_privacy_odbc.beam
lib/erlang/lib/${DISTNAME}/ebin/mod_private.beam
lib/erlang/lib/${DISTNAME}/ebin/mod_private_odbc.beam
lib/erlang/lib/${DISTNAME}/ebin/mod_proxy65.beam
lib/erlang/lib/${DISTNAME}/ebin/mod_proxy65_lib.beam
lib/erlang/lib/${DISTNAME}/ebin/mod_proxy65_service.beam
lib/erlang/lib/${DISTNAME}/ebin/mod_proxy65_sm.beam
lib/erlang/lib/${DISTNAME}/ebin/mod_proxy65_stream.beam
lib/erlang/lib/${DISTNAME}/ebin/mod_pubsub.beam
lib/erlang/lib/${DISTNAME}/ebin/mod_register.beam
lib/erlang/lib/${DISTNAME}/ebin/mod_roster.beam
lib/erlang/lib/${DISTNAME}/ebin/mod_roster_odbc.beam
lib/erlang/lib/${DISTNAME}/ebin/mod_service_log.beam
lib/erlang/lib/${DISTNAME}/ebin/mod_shared_roster.beam
lib/erlang/lib/${DISTNAME}/ebin/mod_stats.beam
lib/erlang/lib/${DISTNAME}/ebin/mod_time.beam
lib/erlang/lib/${DISTNAME}/ebin/mod_vcard.beam
lib/erlang/lib/${DISTNAME}/ebin/mod_vcard_ldap.beam
lib/erlang/lib/${DISTNAME}/ebin/mod_vcard_odbc.beam
lib/erlang/lib/${DISTNAME}/ebin/mod_version.beam
lib/erlang/lib/${DISTNAME}/ebin/node_buddy.beam
lib/erlang/lib/${DISTNAME}/ebin/node_club.beam
lib/erlang/lib/${DISTNAME}/ebin/node_default.beam
lib/erlang/lib/${DISTNAME}/ebin/node_dispatch.beam
lib/erlang/lib/${DISTNAME}/ebin/node_pep.beam
lib/erlang/lib/${DISTNAME}/ebin/node_private.beam
lib/erlang/lib/${DISTNAME}/ebin/node_public.beam
lib/erlang/lib/${DISTNAME}/ebin/nodetree_default.beam
lib/erlang/lib/${DISTNAME}/ebin/nodetree_virtual.beam
lib/erlang/lib/${DISTNAME}/ebin/odbc_queries.beam
lib/erlang/lib/${DISTNAME}/ebin/p1_fsm.beam
lib/erlang/lib/${DISTNAME}/ebin/p1_mnesia.beam
lib/erlang/lib/${DISTNAME}/ebin/ram_file_io_server.beam
lib/erlang/lib/${DISTNAME}/ebin/randoms.beam
lib/erlang/lib/${DISTNAME}/ebin/sha.beam
lib/erlang/lib/${DISTNAME}/ebin/shaper.beam
lib/erlang/lib/${DISTNAME}/ebin/stringprep.beam
lib/erlang/lib/${DISTNAME}/ebin/stringprep_sup.beam
lib/erlang/lib/${DISTNAME}/ebin/tls.beam
lib/erlang/lib/${DISTNAME}/ebin/translate.beam
lib/erlang/lib/${DISTNAME}/ebin/xml.beam
lib/erlang/lib/${DISTNAME}/ebin/xml_stream.beam
lib/erlang/lib/${DISTNAME}/priv/lib/ejabberd_zlib_drv.so
lib/erlang/lib/${DISTNAME}/priv/lib/expat_erl.so
lib/erlang/lib/${DISTNAME}/priv/lib/iconv_erl.so
lib/erlang/lib/${DISTNAME}/priv/lib/stringprep_drv.so
lib/erlang/lib/${DISTNAME}/priv/lib/tls_drv.so
lib/erlang/lib/${DISTNAME}/priv/msgs/ca.msg
lib/erlang/lib/${DISTNAME}/priv/msgs/cs.msg
lib/erlang/lib/${DISTNAME}/priv/msgs/de.msg
lib/erlang/lib/${DISTNAME}/priv/msgs/eo.msg
lib/erlang/lib/${DISTNAME}/priv/msgs/es.msg
lib/erlang/lib/${DISTNAME}/priv/msgs/fr.msg
lib/erlang/lib/${DISTNAME}/priv/msgs/gl.msg
lib/erlang/lib/${DISTNAME}/priv/msgs/it.msg
lib/erlang/lib/${DISTNAME}/priv/msgs/ja.msg
lib/erlang/lib/${DISTNAME}/priv/msgs/nl.msg
lib/erlang/lib/${DISTNAME}/priv/msgs/no.msg
lib/erlang/lib/${DISTNAME}/priv/msgs/pl.msg
lib/erlang/lib/${DISTNAME}/priv/msgs/pt-br.msg
lib/erlang/lib/${DISTNAME}/priv/msgs/pt.msg
lib/erlang/lib/${DISTNAME}/priv/msgs/ru.msg
lib/erlang/lib/${DISTNAME}/priv/msgs/sk.msg
lib/erlang/lib/${DISTNAME}/priv/msgs/sv.msg
lib/erlang/lib/${DISTNAME}/priv/msgs/th.msg
lib/erlang/lib/${DISTNAME}/priv/msgs/tr.msg
lib/erlang/lib/${DISTNAME}/priv/msgs/uk.msg
lib/erlang/lib/${DISTNAME}/priv/msgs/vi.msg
lib/erlang/lib/${DISTNAME}/priv/msgs/wa.msg
lib/erlang/lib/${DISTNAME}/priv/msgs/zh.msg
sbin/ejabberdctl
share/doc/${DISTNAME}/Makefile
share/doc/${DISTNAME}/COPYING
share/doc/${DISTNAME}/dev.html
share/doc/${DISTNAME}/dev.tex
share/doc/${DISTNAME}/disco.png
share/doc/${DISTNAME}/discorus.png
share/doc/${DISTNAME}/features.html
share/doc/${DISTNAME}/features.tex
share/doc/${DISTNAME}/flow.dot
share/doc/${DISTNAME}/guide.html
share/doc/${DISTNAME}/guide.tex
share/doc/${DISTNAME}/guide.pdf
share/doc/${DISTNAME}/introduction.tex
share/doc/${DISTNAME}/logo.png
share/doc/${DISTNAME}/mod_http_bind.tex
share/doc/${DISTNAME}/mod_http_fileserver.tex
share/doc/${DISTNAME}/release_notes_0.9.1.txt
share/doc/${DISTNAME}/release_notes_0.9.8.txt
share/doc/${DISTNAME}/release_notes_0.9.txt
share/doc/${DISTNAME}/release_notes_1.0.0.txt
share/doc/${DISTNAME}/release_notes_1.1.0.txt
share/doc/${DISTNAME}/release_notes_1.1.1.txt
share/doc/${DISTNAME}/release_notes_1.1.2.txt
share/doc/${DISTNAME}/release_notes_1.1.3.txt
share/doc/${DISTNAME}/release_notes_1.1.4.txt
share/doc/${DISTNAME}/release_notes_2.0.0.txt
share/doc/${DISTNAME}/version.tex
share/doc/${DISTNAME}/webadmmain.png
share/doc/${DISTNAME}/webadmmainru.png
share/doc/${DISTNAME}/yozhikheader.png
share/doc/${DISTNAME}/api/Makefile
share/doc/${DISTNAME}/api/process-one.css
share/doc/${DISTNAME}/api/overview.edoc
share/examples/rc.d/ejabberd
@dirrm share/examples/ejabberd
@dirrm share/doc/${DISTNAME}/api
@dirrm share/doc/${DISTNAME}
@dirrm lib/erlang/lib/${DISTNAME}/priv/msgs
@dirrm lib/erlang/lib/${DISTNAME}/priv/lib
@dirrm lib/erlang/lib/${DISTNAME}/priv
@dirrm lib/erlang/lib/${DISTNAME}/ebin
@dirrm lib/erlang/lib/${DISTNAME}
@unexec ${RMDIR} %D/etc/ejabberd 2>/dev/null || ${TRUE}
