# net/ipa/PLIST
@comment $NetBSD: PLIST,v 1.5 2005/05/02 20:34:02 reed Exp $
bin/ipa
bin/ipastat
man/man5/ipa.5
man/man5/ipa.conf.5
man/man8/ipa.8
man/man8/ipastat.8
share/examples/ipa/README
share/examples/ipa/example-1
share/examples/ipa/example-2
share/examples/ipa/ipa.conf.default
share/examples/rc.d/ipa
@dirrm share/examples/ipa
