# sysutils/nabud/PLIST.NetBSD
@comment $NetBSD: PLIST.NetBSD,v 1.1 2023/01/15 21:41:34 thorpej Exp $
share/nabud/rc.conf.d/nabud
share/nabud/rc.d/nabud
