# security/cfs/PLIST
@comment $NetBSD: PLIST,v 1.3 2005/05/02 20:34:04 reed Exp $
bin/cfs_attach
bin/cfs_cat
bin/cfs_detach
bin/cfs_mkdir
bin/cfs_mkkey
bin/cfs_name
bin/cfs_passwd
bin/cfssh
man/man1/cfs_attach.1
man/man1/cfs_detach.1
man/man1/cfs_mkdir.1
man/man1/cfs_mkkey.1
man/man1/cfs_passwd.1
man/man1/cfssh.1
man/man8/cfs_cat.8
man/man8/cfs_name.8
man/man8/cfsd.8
sbin/cfsd
share/doc/cfs/notes.ms
share/examples/rc.d/cfsd
@dirrm share/doc/cfs
