# comms/zaptel-netbsd/PLIST
@comment $NetBSD: PLIST,v 1.4 2009/06/14 18:31:59 joerg Exp $
bin/fxo_test
bin/ztcfg
bin/ztmonitor
bin/zttest
bin/zttool
include/tonezone.h
include/zaptel.h
lib/libtonezone.a
lkm/mod_wcfxo.o
lkm/mod_wcfxs.o
lkm/mod_zaptel.o
lkm/mod_ztdummy.o
sbin/zaptel_postinstall
share/examples/rc.d/ztcfg
share/examples/zaptel/zaptel.conf
