# devel/memcached/PLIST
@comment $NetBSD: PLIST,v 1.7 2011/10/03 11:53:11 adam Exp $
bin/memcached
bin/memcached-tool
include/memcached/protocol_binary.h
man/man1/memcached.1
share/doc/memcached/protocol.txt
share/doc/memcached/readme.txt
share/examples/rc.d/memcached
