# devel/memcached/PLIST
@comment $NetBSD: PLIST,v 1.8 2012/10/04 10:46:48 asau Exp $
bin/memcached
bin/memcached-tool
include/memcached/protocol_binary.h
man/man1/memcached.1
share/doc/memcached/protocol.txt
share/doc/memcached/readme.txt
share/doc/memcached/threads.txt
share/examples/rc.d/memcached
