# www/ruby-puma/PLIST
@comment $NetBSD: PLIST,v 1.17 2023/02/23 15:41:17 taca Exp $
bin/puma${RUBY_SUFFIX}
bin/pumactl${RUBY_SUFFIX}
${GEM_HOME}/cache/${GEM_NAME}.gem
${GEM_EXTSDIR}/gem.build_complete
${GEM_EXTSDIR}/puma/puma_http11.${RUBY_DLEXT}
${GEM_LIBDIR}/History.md
${GEM_LIBDIR}/LICENSE
${GEM_LIBDIR}/README.md
${GEM_LIBDIR}/bin/puma
${GEM_LIBDIR}/bin/puma-wild
${GEM_LIBDIR}/bin/pumactl
${GEM_LIBDIR}/docs/architecture.md
${GEM_LIBDIR}/docs/compile_options.md
${GEM_LIBDIR}/docs/deployment.md
${GEM_LIBDIR}/docs/fork_worker.md
${GEM_LIBDIR}/docs/images/puma-connection-flow-no-reactor.png
${GEM_LIBDIR}/docs/images/puma-connection-flow.png
${GEM_LIBDIR}/docs/images/puma-general-arch.png
${GEM_LIBDIR}/docs/jungle/README.md
${GEM_LIBDIR}/docs/jungle/rc.d/README.md
${GEM_LIBDIR}/docs/jungle/rc.d/puma
${GEM_LIBDIR}/docs/jungle/rc.d/puma.conf
${GEM_LIBDIR}/docs/kubernetes.md
${GEM_LIBDIR}/docs/nginx.md
${GEM_LIBDIR}/docs/plugins.md
${GEM_LIBDIR}/docs/rails_dev_mode.md
${GEM_LIBDIR}/docs/restart.md
${GEM_LIBDIR}/docs/signals.md
${GEM_LIBDIR}/docs/stats.md
${GEM_LIBDIR}/docs/systemd.md
${GEM_LIBDIR}/docs/testing_benchmarks_local_files.md
${GEM_LIBDIR}/docs/testing_test_rackup_ci_files.md
${GEM_LIBDIR}/ext/puma_http11/PumaHttp11Service.java
${GEM_LIBDIR}/ext/puma_http11/ext_help.h
${GEM_LIBDIR}/ext/puma_http11/extconf.rb
${GEM_LIBDIR}/ext/puma_http11/http11_parser.c
${GEM_LIBDIR}/ext/puma_http11/http11_parser.h
${GEM_LIBDIR}/ext/puma_http11/http11_parser.java.rl
${GEM_LIBDIR}/ext/puma_http11/http11_parser.rl
${GEM_LIBDIR}/ext/puma_http11/http11_parser_common.rl
${GEM_LIBDIR}/ext/puma_http11/mini_ssl.c
${GEM_LIBDIR}/ext/puma_http11/no_ssl/PumaHttp11Service.java
${GEM_LIBDIR}/ext/puma_http11/org/jruby/puma/Http11.java
${GEM_LIBDIR}/ext/puma_http11/org/jruby/puma/Http11Parser.java
${GEM_LIBDIR}/ext/puma_http11/org/jruby/puma/MiniSSL.java
${GEM_LIBDIR}/ext/puma_http11/puma_http11.c
${GEM_LIBDIR}/lib/puma.rb
${GEM_LIBDIR}/lib/puma/app/status.rb
${GEM_LIBDIR}/lib/puma/binder.rb
${GEM_LIBDIR}/lib/puma/cli.rb
${GEM_LIBDIR}/lib/puma/client.rb
${GEM_LIBDIR}/lib/puma/cluster.rb
${GEM_LIBDIR}/lib/puma/cluster/worker.rb
${GEM_LIBDIR}/lib/puma/cluster/worker_handle.rb
${GEM_LIBDIR}/lib/puma/commonlogger.rb
${GEM_LIBDIR}/lib/puma/configuration.rb
${GEM_LIBDIR}/lib/puma/const.rb
${GEM_LIBDIR}/lib/puma/control_cli.rb
${GEM_LIBDIR}/lib/puma/detect.rb
${GEM_LIBDIR}/lib/puma/dsl.rb
${GEM_LIBDIR}/lib/puma/error_logger.rb
${GEM_LIBDIR}/lib/puma/events.rb
${GEM_LIBDIR}/lib/puma/io_buffer.rb
${GEM_LIBDIR}/lib/puma/jruby_restart.rb
${GEM_LIBDIR}/lib/puma/json_serialization.rb
${GEM_LIBDIR}/lib/puma/launcher.rb
${GEM_LIBDIR}/lib/puma/launcher/bundle_pruner.rb
${GEM_LIBDIR}/lib/puma/log_writer.rb
${GEM_LIBDIR}/lib/puma/minissl.rb
${GEM_LIBDIR}/lib/puma/minissl/context_builder.rb
${GEM_LIBDIR}/lib/puma/null_io.rb
${GEM_LIBDIR}/lib/puma/plugin.rb
${GEM_LIBDIR}/lib/puma/plugin/systemd.rb
${GEM_LIBDIR}/lib/puma/plugin/tmp_restart.rb
${GEM_LIBDIR}/lib/puma/puma_http11.${RUBY_DLEXT}
${GEM_LIBDIR}/lib/puma/rack/builder.rb
${GEM_LIBDIR}/lib/puma/rack/urlmap.rb
${GEM_LIBDIR}/lib/puma/rack_default.rb
${GEM_LIBDIR}/lib/puma/reactor.rb
${GEM_LIBDIR}/lib/puma/request.rb
${GEM_LIBDIR}/lib/puma/runner.rb
${GEM_LIBDIR}/lib/puma/sd_notify.rb
${GEM_LIBDIR}/lib/puma/server.rb
${GEM_LIBDIR}/lib/puma/single.rb
${GEM_LIBDIR}/lib/puma/state_file.rb
${GEM_LIBDIR}/lib/puma/thread_pool.rb
${GEM_LIBDIR}/lib/puma/util.rb
${GEM_LIBDIR}/lib/rack/handler/puma.rb
${GEM_LIBDIR}/tools/Dockerfile
${GEM_LIBDIR}/tools/trickletest.rb
${GEM_HOME}/specifications/${GEM_NAME}.gemspec
