# inputmethod/skkserv/PLIST
@comment $NetBSD: PLIST,v 1.1.1.1.6.2 2002/06/23 18:48:21 jlam Exp $
bin/skkdic-expr
bin/skkdic-sort
bin/skkdic-count
etc/rc.d/skkserv
libexec/skkserv
share/skk/SKK-JISYO.L
share/skk/README.skkserv
share/skk/README.1st.jisyo-tools
share/skk/README.jisyo-tools
@dirrm share/skk
