# inputmethod/skkserv/PLIST
@comment $NetBSD: PLIST,v 1.3 2005/05/02 20:33:59 reed Exp $
bin/skkdic-expr
bin/skkdic-sort
bin/skkdic-count
libexec/skkserv
share/examples/rc.d/skkserv
share/skk/SKK-JISYO.L
share/skk/README.skkserv
share/skk/README.1st.jisyo-tools
share/skk/README.jisyo-tools
@dirrm share/skk
