# inputmethod/skkserv/PLIST
@comment $NetBSD: PLIST,v 1.1.1.1 2002/05/31 12:59:58 seb Exp $
bin/skkdic-expr
bin/skkdic-sort
bin/skkdic-count
etc/rc.d/skkserv
libexec/skkserv
share/skk/SKK-JISYO.L
share/skk/README.skkserv
share/skk/README.1st.jisyo-tools
share/skk/README.jisyo-tools
@dirrm share/skk
