# mail/dkim-milter/PLIST
@comment $NetBSD: PLIST,v 1.1.1.1 2005/10/27 19:49:07 minskim Exp $
man/man8/dkim-milter.8
sbin/dkim-milter
share/doc/dkim-milter/FEATURES
share/doc/dkim-milter/INSTALL
share/doc/dkim-milter/KNOWNBUGS
share/doc/dkim-milter/README
share/doc/dkim-milter/RELEASE_NOTES
share/examples/rc.d/dkim-milter
@dirrm share/doc/dkim-milter
