# mail/dkim-milter/PLIST
@comment $NetBSD: PLIST,v 1.8 2008/04/20 17:36:15 adrianp Exp $
${PLIST.arlib}lib/libar.a
${PLIST.arlib}man/man3/ar.3
lib/libdkim.a
include/dkim.h
libexec/dkim-filter
bin/dkim-genkey
${PLIST.stats}bin/dkim-stats
bin/dkim-testkey
bin/dkim-testssp
man/man5/dkim-filter.conf.5
man/man8/dkim-filter.8
man/man8/dkim-genkey.8
${PLIST.stats}man/man8/dkim-stats.8
man/man8/dkim-testkey.8
man/man8/dkim-testssp.8
share/doc/dkim-milter/FEATURES
share/doc/dkim-milter/INSTALL
share/doc/dkim-milter/KNOWNBUGS
share/doc/dkim-milter/README
share/doc/dkim-milter/RELEASE_NOTES
share/examples/dkim-milter/dkim-filter.conf.sample
share/doc/dkim-milter/draft-hoffman-dac-vbr-02.txt
share/doc/dkim-milter/draft-ietf-dkim-deployment-01.txt
share/doc/dkim-milter/draft-ietf-dkim-overview-09.txt
share/doc/dkim-milter/draft-ietf-dkim-ssp-03.txt
share/doc/dkim-milter/draft-kucherawy-sender-auth-header-14.txt
share/doc/dkim-milter/libdkim/dkim.html
share/doc/dkim-milter/libdkim/dkim_alg_t.html
share/doc/dkim-milter/libdkim/dkim_body.html
share/doc/dkim-milter/libdkim/dkim_canon_t.html
share/doc/dkim-milter/libdkim/dkim_cbstat.html
share/doc/dkim-milter/libdkim/dkim_close.html
share/doc/dkim-milter/libdkim/dkim_eoh.html
share/doc/dkim-milter/libdkim/dkim_eom.html
share/doc/dkim-milter/libdkim/dkim_flush_cache.html
share/doc/dkim-milter/libdkim/dkim_free.html
share/doc/dkim-milter/libdkim/dkim_get_user_context.html
share/doc/dkim-milter/libdkim/dkim_getcachestats.html
share/doc/dkim-milter/libdkim/dkim_getdomain.html
share/doc/dkim-milter/libdkim/dkim_geterror.html
share/doc/dkim-milter/libdkim/dkim_getmode.html
share/doc/dkim-milter/libdkim/dkim_getpolicystr.html
share/doc/dkim-milter/libdkim/dkim_getpresult.html
share/doc/dkim-milter/libdkim/dkim_getpresultstr.html
share/doc/dkim-milter/libdkim/dkim_getresultstr.html
share/doc/dkim-milter/libdkim/dkim_getsighdr.html
share/doc/dkim-milter/libdkim/dkim_getsiglist.html
share/doc/dkim-milter/libdkim/dkim_getsignature.html
share/doc/dkim-milter/libdkim/dkim_header.html
share/doc/dkim-milter/libdkim/dkim_init.html
share/doc/dkim-milter/libdkim/dkim_lib.html
share/doc/dkim-milter/libdkim/dkim_minbody.html
share/doc/dkim-milter/libdkim/dkim_ohdrs.html
share/doc/dkim-milter/libdkim/dkim_options.html
share/doc/dkim-milter/libdkim/dkim_param_t.html
share/doc/dkim-milter/libdkim/dkim_policy.html
share/doc/dkim-milter/libdkim/dkim_policy_t.html
share/doc/dkim-milter/libdkim/dkim_presult.html
share/doc/dkim-milter/libdkim/dkim_pstate.html
share/doc/dkim-milter/libdkim/dkim_query_t.html
share/doc/dkim-milter/libdkim/dkim_reportinfo.html
share/doc/dkim-milter/libdkim/dkim_set_dns_callback.html
share/doc/dkim-milter/libdkim/dkim_set_final.html
share/doc/dkim-milter/libdkim/dkim_set_key_lookup.html
share/doc/dkim-milter/libdkim/dkim_set_policy_lookup.html
share/doc/dkim-milter/libdkim/dkim_set_prescreen.html
share/doc/dkim-milter/libdkim/dkim_set_signature_handle.html
share/doc/dkim-milter/libdkim/dkim_set_signature_handle_free.html
share/doc/dkim-milter/libdkim/dkim_set_signature_tagvalues.html
share/doc/dkim-milter/libdkim/dkim_set_signer.html
share/doc/dkim-milter/libdkim/dkim_set_user_context.html
share/doc/dkim-milter/libdkim/dkim_sig_getbh.html
share/doc/dkim-milter/libdkim/dkim_sig_getcanonlen.html
share/doc/dkim-milter/libdkim/dkim_sig_getcontext.html
share/doc/dkim-milter/libdkim/dkim_sig_getdomain.html
share/doc/dkim-milter/libdkim/dkim_sig_geterror.html
share/doc/dkim-milter/libdkim/dkim_sig_geterrorstr.html
share/doc/dkim-milter/libdkim/dkim_sig_getflags.html
share/doc/dkim-milter/libdkim/dkim_sig_getkeysize.html
share/doc/dkim-milter/libdkim/dkim_sig_getselector.html
share/doc/dkim-milter/libdkim/dkim_sig_getsignalg.html
share/doc/dkim-milter/libdkim/dkim_sig_getsigntime.html
share/doc/dkim-milter/libdkim/dkim_sig_ignore.html
share/doc/dkim-milter/libdkim/dkim_sig_process.html
share/doc/dkim-milter/libdkim/dkim_sigerror.html
share/doc/dkim-milter/libdkim/dkim_siginfo.html
share/doc/dkim-milter/libdkim/dkim_sigkey_t.html
share/doc/dkim-milter/libdkim/dkim_sign.html
share/doc/dkim-milter/libdkim/dkim_ssl_version.html
share/doc/dkim-milter/libdkim/dkim_stat.html
share/doc/dkim-milter/libdkim/dkim_verify.html
share/doc/dkim-milter/libdkim/index.html
share/doc/dkim-milter/libdkim/dkim_chunk.html
share/doc/dkim-milter/libdkim/dkim_key_syntax.html
share/doc/dkim-milter/libdkim/dkim_policy_syntax.html
share/doc/dkim-milter/libdkim/dkim_sig_syntax.html
share/doc/dkim-milter/libdkim/overview.html
share/doc/dkim-milter/libdkim/rfc2822_mailbox_split.html
share/doc/dkim-milter/rfc4686.txt
share/doc/dkim-milter/rfc4871.txt
share/doc/dkim-milter/rfc5016.txt
share/examples/rc.d/dkim-filter
@dirrm share/doc/dkim-milter/libdkim
@dirrm share/doc/dkim-milter
@dirrm share/examples/dkim-milter
