# security/prelude-lml/PLIST
@comment $NetBSD: PLIST,v 1.7 2007/05/12 10:00:35 shannonjr Exp $
bin/prelude-lml
include/prelude-lml/prelude-lml.h
lib/prelude-lml/debug.la
lib/prelude-lml/pcre.la
sbin/run-prelude-lml
share/examples/prelude-lml/plugins.rules
share/examples/prelude-lml/prelude-lml.conf
share/examples/prelude-lml/ruleset/apc-emu.rules
share/examples/prelude-lml/ruleset/arbor.rules
share/examples/prelude-lml/ruleset/arpwatch.rules
share/examples/prelude-lml/ruleset/checkpoint.rules
share/examples/prelude-lml/ruleset/cisco-router.rules
share/examples/prelude-lml/ruleset/cisco-vpn.rules
share/examples/prelude-lml/ruleset/clamav.rules
share/examples/prelude-lml/ruleset/dell-om.rules
share/examples/prelude-lml/ruleset/f5-bigip.rules
share/examples/prelude-lml/ruleset/grsecurity.rules
share/examples/prelude-lml/ruleset/honeyd.rules
share/examples/prelude-lml/ruleset/httpd.rules
share/examples/prelude-lml/ruleset/ipchains.rules
share/examples/prelude-lml/ruleset/ipfw.rules
share/examples/prelude-lml/ruleset/ipso.rules
share/examples/prelude-lml/ruleset/linksys-wap11.rules
share/examples/prelude-lml/ruleset/modsecurity.rules
share/examples/prelude-lml/ruleset/ms-sql.rules
share/examples/prelude-lml/ruleset/nagios.rules
share/examples/prelude-lml/ruleset/navce.rules
share/examples/prelude-lml/ruleset/netapp-ontap.rules
share/examples/prelude-lml/ruleset/netfilter.rules
share/examples/prelude-lml/ruleset/netscreen.rules
share/examples/prelude-lml/ruleset/ntsyslog.rules
share/examples/prelude-lml/ruleset/openhostapd.rules
share/examples/prelude-lml/ruleset/pam.rules
share/examples/prelude-lml/ruleset/pcanywhere.rules
share/examples/prelude-lml/ruleset/pcre.rules
share/examples/prelude-lml/ruleset/portsentry.rules
share/examples/prelude-lml/ruleset/postfix.rules
share/examples/prelude-lml/ruleset/proftpd.rules
share/examples/prelude-lml/ruleset/qpopper.rules
share/examples/prelude-lml/ruleset/selinux.rules
share/examples/prelude-lml/ruleset/sendmail.rules
share/examples/prelude-lml/ruleset/shadow-utils.rules
share/examples/prelude-lml/ruleset/single.rules
share/examples/prelude-lml/ruleset/sonicwall.rules
share/examples/prelude-lml/ruleset/spamassassin.rules
share/examples/prelude-lml/ruleset/squid.rules
share/examples/prelude-lml/ruleset/ssh.rules
share/examples/prelude-lml/ruleset/sudo.rules
share/examples/prelude-lml/ruleset/tripwire.rules
share/examples/prelude-lml/ruleset/vigor.rules
share/examples/prelude-lml/ruleset/vpopmail.rules
share/examples/prelude-lml/ruleset/webmin.rules
share/examples/prelude-lml/ruleset/wu-ftp.rules
share/examples/prelude-lml/ruleset/zywall.rules
share/examples/prelude-lml/ruleset/zyxel.rules
share/examples/prelude-lml/ruleset/cacti-thold.rules
share/examples/prelude-lml/ruleset/cisco-asa.rules
share/examples/prelude-lml/ruleset/cisco-css.rules
share/examples/prelude-lml/ruleset/cisco-ips.rules
share/examples/prelude-lml/ruleset/exim.rules
share/examples/prelude-lml/ruleset/ms-cluster.rules
share/examples/prelude-lml/ruleset/bonding.rules
share/examples/prelude-lml/ruleset/cisco-common.rules
share/examples/rc.d/preludelml
@dirrm share/examples/prelude-lml/ruleset
@dirrm share/examples/prelude-lml
@dirrm lib/prelude-lml
@dirrm include/prelude-lml
