# mail/qmail-run/PLIST
@comment $NetBSD: PLIST,v 1.4 2005/05/03 15:37:28 wiz Exp $
bin/qmail-procmail
bin/qmail-qfilter-queue
bin/qmail-qread-client
share/doc/qmail-run/README.pkgsrc
share/examples/qmail-run/mailer.conf
share/examples/rc.d/qmail
share/examples/rc.d/qmailpop3d
share/examples/rc.d/qmailqread
share/examples/rc.d/qmailsend
share/examples/rc.d/qmailsmtpd
@dirrm share/examples/qmail-run
@dirrm share/doc/qmail-run
