# net/ra-rtsp-proxy/PLIST
@comment $NetBSD: PLIST,v 1.3 2009/06/14 18:09:41 joerg Exp $
man/man8/rtspd.8
sbin/rtspd
share/doc/ra-rtsp-proxy/README
share/doc/ra-rtsp-proxy/license.txt
share/examples/rc.d/rtspd
