# net/dynipclient/PLIST
@comment $NetBSD: PLIST,v 1.5 2005/05/02 20:34:02 reed Exp $
bin/dynipadmin
bin/dynipclient
man/man1/dynipadmin.1
man/man1/dynipclient.1
share/examples/dynipclient/dynip.cfg
share/examples/rc.d/dynipclient
@dirrm share/examples/dynipclient
