# net/dynipclient/PLIST
@comment $NetBSD: PLIST,v 1.3 2003/07/30 12:01:26 jmmv Exp $
bin/dynipadmin
bin/dynipclient
etc/rc.d/dynipclient
man/man1/dynipadmin.1
man/man1/dynipclient.1
share/examples/dynipclient/dynip.cfg
@dirrm share/examples/dynipclient
