# net/delegate/PLIST
@comment $NetBSD: PLIST,v 1.4 2005/05/02 20:34:01 reed Exp $
sbin/delegated
share/doc/delegate/manual.html
share/doc/delegate/tutorial-en.html
share/doc/delegate/tutorial-jp.html
share/examples/rc.d/delegated
@dirrm share/doc/delegate
