# net/delegate/PLIST
@comment $NetBSD: PLIST,v 1.2 2003/07/30 11:37:02 jmmv Exp $
sbin/delegated
etc/rc.d/delegated
share/doc/delegate/manual.html
share/doc/delegate/tutorial-en.html
share/doc/delegate/tutorial-jp.html
@dirrm share/doc/delegate
