# net/haproxy/PLIST
@comment $NetBSD: PLIST,v 1.3 2014/02/20 00:14:22 rodent Exp $
man/man1/haproxy.1
sbin/haproxy
share/doc/haproxy/architecture.txt
share/doc/haproxy/configuration.txt
share/doc/haproxy/haproxy-en.txt
share/doc/haproxy/haproxy-fr.txt
share/examples/haproxy/haproxy.cfg
share/examples/rc.d/haproxy
