# net/haproxy/PLIST
@comment $NetBSD: PLIST,v 1.2 2011/12/16 12:14:12 fhajny Exp $
man/man1/haproxy.1
sbin/haproxy
share/doc/haproxy/architecture.txt
share/doc/haproxy/configuration.txt
share/doc/haproxy/haproxy-en.txt
share/doc/haproxy/haproxy-fr.txt
share/examples/haproxy/haproxy.cfg
${PLIST.rcd}share/examples/rc.d/haproxy
