# net/vnstat/PLIST
@comment $NetBSD: PLIST,v 1.2 2009/11/26 14:11:44 joerg Exp $
bin/vnstat
bin/vnstatd
${PLIST.gd}bin/vnstati
man/man1/vnstat.1
man/man1/vnstatd.1
${PLIST.gd}man/man1/vnstati.1
man/man5/vnstat.conf.5
share/examples/vnstat.conf
share/examples/rc.d/vnstatd
