# sysutils/storage-manager/files/dpteng
#!/bin/sh
#
# $NetBSD: dpteng,v 1.1.1.1 2001/09/26 07:28:49 ad Exp $
#

# PROVIDE: dpteng
# REQUIRE: DAEMON

. /etc/rc.subr

name="dpteng"
rcvar=$name
command="@PREFIX@/lib/dpt/${name}"
pidfile="/tmp/dpteng.LOCK"
start_cmd="dpteng_start"
stop_cmd="dpteng_stop"

dpteng_start()
{

	echo Starting $name.
	$command 2> /dev/null &
}

dpteng_stop()
{

	echo Stopping $name.
	kill -9 `cat $pidfile`
}

load_rc_config $name
run_rc_command "$1" 
