# mail/fetchmail/PLIST
@comment $NetBSD: PLIST,v 1.9 2005/05/02 20:34:00 reed Exp $
bin/fetchmail
man/man1/fetchmail.1
share/doc/fetchmail/README
share/doc/fetchmail/COPYING
share/doc/fetchmail/FAQ
share/doc/fetchmail/NOTES
share/doc/fetchmail/FEATURES
share/doc/fetchmail/fetchmail-FAQ.html
share/doc/fetchmail/design-notes.html
share/doc/fetchmail/fetchmail-features.html
${PKGLOCALEDIR}/locale/ca/LC_MESSAGES/fetchmail.mo
${PKGLOCALEDIR}/locale/cs/LC_MESSAGES/fetchmail.mo
${PKGLOCALEDIR}/locale/da/LC_MESSAGES/fetchmail.mo
${PKGLOCALEDIR}/locale/de/LC_MESSAGES/fetchmail.mo
${PKGLOCALEDIR}/locale/el/LC_MESSAGES/fetchmail.mo
${PKGLOCALEDIR}/locale/es/LC_MESSAGES/fetchmail.mo
${PKGLOCALEDIR}/locale/fr/LC_MESSAGES/fetchmail.mo
${PKGLOCALEDIR}/locale/gl/LC_MESSAGES/fetchmail.mo
${PKGLOCALEDIR}/locale/ja/LC_MESSAGES/fetchmail.mo
${PKGLOCALEDIR}/locale/pl/LC_MESSAGES/fetchmail.mo
${PKGLOCALEDIR}/locale/pt_BR/LC_MESSAGES/fetchmail.mo
${PKGLOCALEDIR}/locale/sk/LC_MESSAGES/fetchmail.mo
${PKGLOCALEDIR}/locale/tr/LC_MESSAGES/fetchmail.mo
share/examples/rc.d/fetchmail
@dirrm share/doc/fetchmail
