# mail/fetchmail/PLIST
@comment $NetBSD: PLIST,v 1.10.6.1 2007/01/17 17:52:01 salo Exp $
bin/fetchmail
man/man1/fetchmail.1
share/doc/fetchmail/README
share/doc/fetchmail/COPYING
share/doc/fetchmail/FAQ
share/doc/fetchmail/NOTES
share/doc/fetchmail/FEATURES
share/doc/fetchmail/fetchmail-FAQ.html
share/doc/fetchmail/design-notes.html
share/doc/fetchmail/fetchmail-features.html
share/locale/ca/LC_MESSAGES/fetchmail.mo
share/locale/cs/LC_MESSAGES/fetchmail.mo
share/locale/de/LC_MESSAGES/fetchmail.mo
share/locale/en_GB/LC_MESSAGES/fetchmail.mo
share/locale/es/LC_MESSAGES/fetchmail.mo
share/locale/fr/LC_MESSAGES/fetchmail.mo
share/locale/ja/LC_MESSAGES/fetchmail.mo
share/locale/pl/LC_MESSAGES/fetchmail.mo
share/locale/ru/LC_MESSAGES/fetchmail.mo
share/locale/sq/LC_MESSAGES/fetchmail.mo
share/locale/vi/LC_MESSAGES/fetchmail.mo
share/examples/rc.d/fetchmail
@dirrm share/doc/fetchmail
