# net/pxe/PLIST
@comment $NetBSD: PLIST,v 1.4 2009/06/14 18:09:39 joerg Exp $
sbin/pxe
share/examples/pxe/pxe.conf
share/examples/rc.d/pxe
