# filesystems/glusterfs/PLIST
@comment $NetBSD: PLIST,v 1.14 2013/04/06 04:03:36 rodent Exp $
etc/glusterfs/glusterd.vol.sample
lib/libglusterfs.la
lib/libgfrpc.la
lib/libgfxdr.la
lib/glusterfs/3.3.1/auth/addr.la
lib/glusterfs/3.3.1/auth/login.la
lib/glusterfs/3.3.1/rpc-transport/socket.la
lib/glusterfs/3.3.1/xlator/cluster/afr.la
lib/glusterfs/3.3.1/xlator/cluster/dht.la
lib/glusterfs/3.3.1/xlator/cluster/distribute.so
lib/glusterfs/3.3.1/xlator/cluster/nufa.la
lib/glusterfs/3.3.1/xlator/cluster/pump.la
lib/glusterfs/3.3.1/xlator/cluster/replicate.so
lib/glusterfs/3.3.1/xlator/cluster/stripe.la
lib/glusterfs/3.3.1/xlator/cluster/switch.la
lib/glusterfs/3.3.1/xlator/debug/error-gen.la
lib/glusterfs/3.3.1/xlator/debug/io-stats.la
lib/glusterfs/3.3.1/xlator/debug/trace.la
lib/glusterfs/3.3.1/xlator/encryption/rot-13.la
lib/glusterfs/3.3.1/xlator/features/access-control.so
lib/glusterfs/3.3.1/xlator/features/index.la
lib/glusterfs/3.3.1/xlator/features/locks.la
lib/glusterfs/3.3.1/xlator/features/mac-compat.la
lib/glusterfs/3.3.1/xlator/features/marker.la
lib/glusterfs/3.3.1/xlator/features/posix-locks.so
lib/glusterfs/3.3.1/xlator/features/quiesce.la
lib/glusterfs/3.3.1/xlator/features/quota.la
lib/glusterfs/3.3.1/xlator/features/read-only.la
lib/glusterfs/3.3.1/xlator/features/worm.la
lib/glusterfs/3.3.1/xlator/mgmt/glusterd.la
lib/glusterfs/3.3.1/xlator/mount/fuse.la
lib/glusterfs/3.3.1/xlator/nfs/server.la
lib/glusterfs/3.3.1/xlator/performance/io-cache.la
lib/glusterfs/3.3.1/xlator/performance/io-threads.la
lib/glusterfs/3.3.1/xlator/performance/md-cache.la
lib/glusterfs/3.3.1/xlator/performance/quick-read.la
lib/glusterfs/3.3.1/xlator/performance/read-ahead.la
lib/glusterfs/3.3.1/xlator/performance/stat-prefetch.so
lib/glusterfs/3.3.1/xlator/performance/write-behind.la
lib/glusterfs/3.3.1/xlator/protocol/client.la
lib/glusterfs/3.3.1/xlator/protocol/server.la
lib/glusterfs/3.3.1/xlator/storage/posix.la
lib/glusterfs/3.3.1/xlator/system/posix-acl.la
lib/glusterfs/3.3.1/xlator/testing/performance/symlink-cache.la
man/man8/gluster.8
sbin/glusterd
sbin/glusterfs
sbin/glusterfsd
sbin/gluster
sbin/mount.glusterfs
share/doc/glusterfs/benchmarking/rdd.c
share/doc/glusterfs/benchmarking/glfs-bm.c
share/doc/glusterfs/benchmarking/README
share/doc/glusterfs/benchmarking/launch-script.sh
share/doc/glusterfs/benchmarking/local-script.sh
share/doc/glusterfs/glusterfs-mode.el
share/doc/glusterfs/glusterfs.vim
share/examples/rc.d/glusterfsd
share/examples/rc.d/glusterd
@exec test -d ${PKG_SYSCONFDIR}/glusterd -a ! -d ${VARBASE}/lib/glusterd && ${MKDIR} -p ${VARBASE}/lib/ && mv ${PKG_SYSCONFDIR}/glusterd  ${VARBASE}/lib/glusterd && ${LN} -s ${VARBASE}/lib/glusterd  ${PKG_SYSCONFDIR}/glusterd
