# filesystems/glusterfs/PLIST
@comment $NetBSD: PLIST,v 1.35 2016/06/16 04:01:02 manu Exp $
${PYSITELIB}/gluster/__init__.py
${PYSITELIB}/gluster/__init__.pyc
${PYSITELIB}/gluster/__init__.pyo
${PYSITELIB}/gluster/glupy/__init__.py
${PYSITELIB}/gluster/glupy/__init__.pyc
${PYSITELIB}/gluster/glupy/__init__.pyo
bin/glusterfind
etc/ganesha/ganesha-ha.conf.sample
etc/glusterfs/gluster-rsyslog-5.8.conf
etc/glusterfs/gluster-rsyslog-7.2.conf
etc/glusterfs/glusterd.vol.sample
etc/glusterfs/glusterfs-georep-logrotate
etc/glusterfs/glusterfs-logrotate
etc/glusterfs/group-virt.example
etc/glusterfs/logger.conf.example
include/glusterfs/api/glfs-handles.h
include/glusterfs/api/glfs.h
include/glusterfs/byte-order.h
include/glusterfs/call-stub.h
include/glusterfs/checksum.h
include/glusterfs/circ-buff.h
include/glusterfs/client_t.h
include/glusterfs/cluster-syncop.h
include/glusterfs/common-utils.h
include/glusterfs/compat-errno.h
include/glusterfs/compat-uuid.h
include/glusterfs/compat.h
include/glusterfs/ctr-messages.h
include/glusterfs/daemon.h
include/glusterfs/default-args.h
include/glusterfs/defaults.h
include/glusterfs/dict.h
include/glusterfs/event-history.h
include/glusterfs/event.h
include/glusterfs/fd-lk.h
include/glusterfs/fd.h
include/glusterfs/gf-dirent.h
include/glusterfs/gfchangelog/changelog.h
include/glusterfs/gfdb/gfdb_data_store.h
include/glusterfs/gfdb/gfdb_data_store_helper.h
include/glusterfs/gfdb/gfdb_data_store_types.h
include/glusterfs/gfdb/gfdb_mem-types.h
include/glusterfs/gfdb/gfdb_sqlite3.h
include/glusterfs/gfdb/gfdb_sqlite3_helper.h
include/glusterfs/gidcache.h
include/glusterfs/glfs-message-id.h
include/glusterfs/globals.h
include/glusterfs/glusterfs-acl.h
include/glusterfs/glusterfs-fops.h
include/glusterfs/glusterfs.h
include/glusterfs/graph-utils.h
include/glusterfs/hashfn.h
include/glusterfs/iatt.h
include/glusterfs/inode.h
include/glusterfs/iobuf.h
include/glusterfs/latency.h
include/glusterfs/libglusterfs-messages.h
include/glusterfs/list.h
include/glusterfs/lkowner.h
include/glusterfs/locking.h
include/glusterfs/logging.h
include/glusterfs/lvm-defaults.h
include/glusterfs/mem-pool.h
include/glusterfs/mem-types.h
include/glusterfs/options.h
include/glusterfs/parse-utils.h
include/glusterfs/quota-common-utils.h
include/glusterfs/rbthash.h
include/glusterfs/refcount.h
include/glusterfs/revision.h
include/glusterfs/rot-buffs.h
include/glusterfs/rpc/acl3-xdr.h
include/glusterfs/rpc/changelog-xdr.h
include/glusterfs/rpc/cli1-xdr.h
include/glusterfs/rpc/glusterd1-xdr.h
include/glusterfs/rpc/glusterfs-fops.h
include/glusterfs/rpc/glusterfs3-xdr.h
include/glusterfs/rpc/glusterfs3.h
include/glusterfs/rpc/mount3udp.h
include/glusterfs/rpc/msg-nfs3.h
include/glusterfs/rpc/nlm4-xdr.h
include/glusterfs/rpc/nsm-xdr.h
include/glusterfs/rpc/portmap-xdr.h
include/glusterfs/rpc/protocol-common.h
include/glusterfs/rpc/rpc-clnt-ping.h
include/glusterfs/rpc/rpc-clnt.h
include/glusterfs/rpc/rpc-common-xdr.h
include/glusterfs/rpc/rpc-drc.h
include/glusterfs/rpc/rpc-lib-messages.h
include/glusterfs/rpc/rpc-transport.h
include/glusterfs/rpc/rpcsvc-common.h
include/glusterfs/rpc/rpcsvc.h
include/glusterfs/rpc/xdr-common.h
include/glusterfs/rpc/xdr-generic.h
include/glusterfs/rpc/xdr-nfs3.h
include/glusterfs/rpc/xdr-rpc.h
include/glusterfs/rpc/xdr-rpcclnt.h
include/glusterfs/run.h
include/glusterfs/server/authenticate.h
include/glusterfs/server/server-common.h
include/glusterfs/server/server-helpers.h
include/glusterfs/server/server-mem-types.h
include/glusterfs/server/server-messages.h
include/glusterfs/server/server.h
include/glusterfs/stack.h
include/glusterfs/statedump.h
include/glusterfs/store.h
include/glusterfs/strfd.h
include/glusterfs/syncop-utils.h
include/glusterfs/syncop.h
include/glusterfs/syscall.h
include/glusterfs/template-component-messages.h
include/glusterfs/timer.h
include/glusterfs/timespec.h
include/glusterfs/trie.h
include/glusterfs/tw.h
include/glusterfs/upcall-utils.h
include/glusterfs/xlator.h
include/glusterfs/y.tab.h
lib/glusterfs/${GLUSTERFS_VERSION}/auth/addr.la
lib/glusterfs/${GLUSTERFS_VERSION}/auth/login.la
lib/glusterfs/${GLUSTERFS_VERSION}/rpc-transport/socket.la
lib/glusterfs/${GLUSTERFS_VERSION}/xlator/cluster/afr.la
lib/glusterfs/${GLUSTERFS_VERSION}/xlator/cluster/dht.la
lib/glusterfs/${GLUSTERFS_VERSION}/xlator/cluster/disperse.so
lib/glusterfs/${GLUSTERFS_VERSION}/xlator/cluster/distribute.so
lib/glusterfs/${GLUSTERFS_VERSION}/xlator/cluster/ec.la
lib/glusterfs/${GLUSTERFS_VERSION}/xlator/cluster/nufa.la
lib/glusterfs/${GLUSTERFS_VERSION}/xlator/cluster/pump.la
lib/glusterfs/${GLUSTERFS_VERSION}/xlator/cluster/replicate.so
lib/glusterfs/${GLUSTERFS_VERSION}/xlator/cluster/stripe.la
lib/glusterfs/${GLUSTERFS_VERSION}/xlator/cluster/switch.la
lib/glusterfs/${GLUSTERFS_VERSION}/xlator/cluster/tier.la
lib/glusterfs/${GLUSTERFS_VERSION}/xlator/debug/error-gen.la
lib/glusterfs/${GLUSTERFS_VERSION}/xlator/debug/io-stats.la
lib/glusterfs/${GLUSTERFS_VERSION}/xlator/debug/trace.la
lib/glusterfs/${GLUSTERFS_VERSION}/xlator/encryption/crypt.la
lib/glusterfs/${GLUSTERFS_VERSION}/xlator/encryption/rot-13.la
lib/glusterfs/${GLUSTERFS_VERSION}/xlator/features/access-control.so
lib/glusterfs/${GLUSTERFS_VERSION}/xlator/features/arbiter.la
lib/glusterfs/${GLUSTERFS_VERSION}/xlator/features/barrier.la
lib/glusterfs/${GLUSTERFS_VERSION}/xlator/features/bit-rot.la
lib/glusterfs/${GLUSTERFS_VERSION}/xlator/features/bitrot-stub.la
lib/glusterfs/${GLUSTERFS_VERSION}/xlator/features/cdc.la
lib/glusterfs/${GLUSTERFS_VERSION}/xlator/features/changelog.la
lib/glusterfs/${GLUSTERFS_VERSION}/xlator/features/changetimerecorder.la
lib/glusterfs/${GLUSTERFS_VERSION}/xlator/features/ganesha.la
lib/glusterfs/${GLUSTERFS_VERSION}/xlator/features/gfid-access.la
lib/glusterfs/${GLUSTERFS_VERSION}/xlator/features/glupy.la
lib/glusterfs/${GLUSTERFS_VERSION}/xlator/features/glupy/debug-trace.py
lib/glusterfs/${GLUSTERFS_VERSION}/xlator/features/glupy/debug-trace.pyc
lib/glusterfs/${GLUSTERFS_VERSION}/xlator/features/glupy/debug-trace.pyo
lib/glusterfs/${GLUSTERFS_VERSION}/xlator/features/glupy/helloworld.py
lib/glusterfs/${GLUSTERFS_VERSION}/xlator/features/glupy/helloworld.pyc
lib/glusterfs/${GLUSTERFS_VERSION}/xlator/features/glupy/helloworld.pyo
lib/glusterfs/${GLUSTERFS_VERSION}/xlator/features/glupy/negative.py
lib/glusterfs/${GLUSTERFS_VERSION}/xlator/features/glupy/negative.pyc
lib/glusterfs/${GLUSTERFS_VERSION}/xlator/features/glupy/negative.pyo
lib/glusterfs/${GLUSTERFS_VERSION}/xlator/features/index.la
lib/glusterfs/${GLUSTERFS_VERSION}/xlator/features/leases.a
lib/glusterfs/${GLUSTERFS_VERSION}/xlator/features/leases.la
lib/glusterfs/${GLUSTERFS_VERSION}/xlator/features/leases.so
lib/glusterfs/${GLUSTERFS_VERSION}/xlator/features/locks.la
lib/glusterfs/${GLUSTERFS_VERSION}/xlator/features/mac-compat.la
lib/glusterfs/${GLUSTERFS_VERSION}/xlator/features/marker.la
lib/glusterfs/${GLUSTERFS_VERSION}/xlator/features/posix-locks.so
lib/glusterfs/${GLUSTERFS_VERSION}/xlator/features/prot_client.la
lib/glusterfs/${GLUSTERFS_VERSION}/xlator/features/prot_dht.la
lib/glusterfs/${GLUSTERFS_VERSION}/xlator/features/prot_server.la
lib/glusterfs/${GLUSTERFS_VERSION}/xlator/features/quiesce.la
lib/glusterfs/${GLUSTERFS_VERSION}/xlator/features/quota.la
lib/glusterfs/${GLUSTERFS_VERSION}/xlator/features/quotad.la
lib/glusterfs/${GLUSTERFS_VERSION}/xlator/features/read-only.la
lib/glusterfs/${GLUSTERFS_VERSION}/xlator/features/shard.la
lib/glusterfs/${GLUSTERFS_VERSION}/xlator/features/snapview-client.la
lib/glusterfs/${GLUSTERFS_VERSION}/xlator/features/snapview-server.la
lib/glusterfs/${GLUSTERFS_VERSION}/xlator/features/trash.la
lib/glusterfs/${GLUSTERFS_VERSION}/xlator/features/upcall.la
lib/glusterfs/${GLUSTERFS_VERSION}/xlator/features/worm.la
lib/glusterfs/${GLUSTERFS_VERSION}/xlator/meta.la
lib/glusterfs/${GLUSTERFS_VERSION}/xlator/mgmt/glusterd.la
lib/glusterfs/${GLUSTERFS_VERSION}/xlator/mount/api.la
lib/glusterfs/${GLUSTERFS_VERSION}/xlator/mount/fuse.la
lib/glusterfs/${GLUSTERFS_VERSION}/xlator/nfs/server.la
lib/glusterfs/${GLUSTERFS_VERSION}/xlator/performance/decompounder.a
lib/glusterfs/${GLUSTERFS_VERSION}/xlator/performance/decompounder.la
lib/glusterfs/${GLUSTERFS_VERSION}/xlator/performance/decompounder.so
lib/glusterfs/${GLUSTERFS_VERSION}/xlator/performance/io-cache.la
lib/glusterfs/${GLUSTERFS_VERSION}/xlator/performance/io-threads.la
lib/glusterfs/${GLUSTERFS_VERSION}/xlator/performance/md-cache.la
lib/glusterfs/${GLUSTERFS_VERSION}/xlator/performance/open-behind.la
lib/glusterfs/${GLUSTERFS_VERSION}/xlator/performance/quick-read.la
lib/glusterfs/${GLUSTERFS_VERSION}/xlator/performance/read-ahead.la
lib/glusterfs/${GLUSTERFS_VERSION}/xlator/performance/readdir-ahead.la
lib/glusterfs/${GLUSTERFS_VERSION}/xlator/performance/stat-prefetch.so
lib/glusterfs/${GLUSTERFS_VERSION}/xlator/performance/write-behind.la
lib/glusterfs/${GLUSTERFS_VERSION}/xlator/protocol/client.la
lib/glusterfs/${GLUSTERFS_VERSION}/xlator/protocol/server.la
lib/glusterfs/${GLUSTERFS_VERSION}/xlator/storage/posix.la
lib/glusterfs/${GLUSTERFS_VERSION}/xlator/system/posix-acl.la
lib/glusterfs/${GLUSTERFS_VERSION}/xlator/testing/features/template.la
lib/glusterfs/${GLUSTERFS_VERSION}/xlator/testing/performance/symlink-cache.la
lib/libgfapi.la
lib/libgfchangelog.la
lib/libgfdb.la
lib/libgfrpc.la
lib/libgfxdr.la
lib/libglusterfs.la
lib/ocf/resource.d/glusterfs/glusterd
lib/ocf/resource.d/glusterfs/volume
lib/ocf/resource.d/heartbeat/ganesha_grace
lib/ocf/resource.d/heartbeat/ganesha_mon
lib/ocf/resource.d/heartbeat/ganesha_nfsd
lib/pkgconfig/glusterfs-api.pc
lib/pkgconfig/libgfchangelog.pc
lib/pkgconfig/libgfdb.pc
libexec/ganesha/copy-export-ganesha.sh
libexec/ganesha/create-export-ganesha.sh
libexec/ganesha/dbus-send.sh
libexec/ganesha/ganesha-ha.sh
libexec/ganesha/generate-epoch.py
libexec/glusterfs/gfind_missing_files/gcrawler
libexec/glusterfs/gfind_missing_files/gfid_to_path.py
libexec/glusterfs/gfind_missing_files/gfid_to_path.sh
libexec/glusterfs/gfind_missing_files/gfind_missing_files.sh
libexec/glusterfs/glusterfind/S57glusterfind-delete-post.py
libexec/glusterfs/glusterfind/__init__.py
libexec/glusterfs/glusterfind/__init__.pyc
libexec/glusterfs/glusterfind/__init__.pyo
libexec/glusterfs/glusterfind/brickfind.py
libexec/glusterfs/glusterfind/changelog.py
libexec/glusterfs/glusterfind/changelogdata.py
libexec/glusterfs/glusterfind/changelogdata.pyc
libexec/glusterfs/glusterfind/changelogdata.pyo
libexec/glusterfs/glusterfind/conf.py
libexec/glusterfs/glusterfind/conf.pyc
libexec/glusterfs/glusterfind/conf.pyo
libexec/glusterfs/glusterfind/libgfchangelog.py
libexec/glusterfs/glusterfind/libgfchangelog.pyc
libexec/glusterfs/glusterfind/libgfchangelog.pyo
libexec/glusterfs/glusterfind/main.py
libexec/glusterfs/glusterfind/main.pyc
libexec/glusterfs/glusterfind/main.pyo
libexec/glusterfs/glusterfind/nodeagent.py
libexec/glusterfs/glusterfind/tool.conf
libexec/glusterfs/glusterfind/utils.py
libexec/glusterfs/glusterfind/utils.pyc
libexec/glusterfs/glusterfind/utils.pyo
libexec/glusterfs/gsyncd
libexec/glusterfs/gverify.sh
libexec/glusterfs/peer_add_secret_pub
libexec/glusterfs/peer_gsec_create
libexec/glusterfs/peer_mountbroker
libexec/glusterfs/python/syncdaemon/README.md
libexec/glusterfs/python/syncdaemon/__init__.py
libexec/glusterfs/python/syncdaemon/__init__.pyc
libexec/glusterfs/python/syncdaemon/__init__.pyo
libexec/glusterfs/python/syncdaemon/changelogagent.py
libexec/glusterfs/python/syncdaemon/changelogagent.pyc
libexec/glusterfs/python/syncdaemon/changelogagent.pyo
libexec/glusterfs/python/syncdaemon/configinterface.py
libexec/glusterfs/python/syncdaemon/configinterface.pyc
libexec/glusterfs/python/syncdaemon/configinterface.pyo
libexec/glusterfs/python/syncdaemon/gconf.py
libexec/glusterfs/python/syncdaemon/gconf.pyc
libexec/glusterfs/python/syncdaemon/gconf.pyo
libexec/glusterfs/python/syncdaemon/gsyncd.py
libexec/glusterfs/python/syncdaemon/gsyncd.pyc
libexec/glusterfs/python/syncdaemon/gsyncd.pyo
libexec/glusterfs/python/syncdaemon/gsyncdstatus.py
libexec/glusterfs/python/syncdaemon/gsyncdstatus.pyc
libexec/glusterfs/python/syncdaemon/gsyncdstatus.pyo
libexec/glusterfs/python/syncdaemon/ipaddr.py
libexec/glusterfs/python/syncdaemon/ipaddr.pyc
libexec/glusterfs/python/syncdaemon/ipaddr.pyo
libexec/glusterfs/python/syncdaemon/libcxattr.py
libexec/glusterfs/python/syncdaemon/libcxattr.pyc
libexec/glusterfs/python/syncdaemon/libcxattr.pyo
libexec/glusterfs/python/syncdaemon/libgfchangelog.py
libexec/glusterfs/python/syncdaemon/libgfchangelog.pyc
libexec/glusterfs/python/syncdaemon/libgfchangelog.pyo
libexec/glusterfs/python/syncdaemon/master.py
libexec/glusterfs/python/syncdaemon/master.pyc
libexec/glusterfs/python/syncdaemon/master.pyo
libexec/glusterfs/python/syncdaemon/monitor.py
libexec/glusterfs/python/syncdaemon/monitor.pyc
libexec/glusterfs/python/syncdaemon/monitor.pyo
libexec/glusterfs/python/syncdaemon/repce.py
libexec/glusterfs/python/syncdaemon/repce.pyc
libexec/glusterfs/python/syncdaemon/repce.pyo
libexec/glusterfs/python/syncdaemon/resource.py
libexec/glusterfs/python/syncdaemon/resource.pyc
libexec/glusterfs/python/syncdaemon/resource.pyo
libexec/glusterfs/python/syncdaemon/syncdutils.py
libexec/glusterfs/python/syncdaemon/syncdutils.pyc
libexec/glusterfs/python/syncdaemon/syncdutils.pyo
libexec/glusterfs/set_geo_rep_pem_keys.sh
man/man8/gluster.8
man/man8/glusterd.8
man/man8/glusterfs.8
man/man8/glusterfsd.8
man/man8/mount.glusterfs.8
sbin/gcron.py
sbin/gfind_missing_files
sbin/glfsheal
sbin/gluster
sbin/glusterd
sbin/glusterfs
sbin/glusterfsd
sbin/mount_glusterfs
sbin/snap_scheduler.py
sbin/umountd
share/doc/glusterfs/benchmarking/README
share/doc/glusterfs/benchmarking/glfs-bm.c
share/doc/glusterfs/benchmarking/launch-script.sh
share/doc/glusterfs/benchmarking/local-script.sh
share/doc/glusterfs/benchmarking/rdd.c
share/doc/glusterfs/glusterfs-mode.el
share/doc/glusterfs/glusterfs.vim
share/examples/rc.d/glusterd
share/glusterfs/scripts/generate-gfid-file.sh
share/glusterfs/scripts/get-gfid.sh
share/glusterfs/scripts/gsync-sync-gfid
share/glusterfs/scripts/gsync-upgrade.sh
share/glusterfs/scripts/post-upgrade-script-for-quota.sh
share/glusterfs/scripts/pre-upgrade-script-for-quota.sh
share/glusterfs/scripts/schedule_georep.py
share/glusterfs/scripts/slave-upgrade.sh
share/glusterfs/scripts/stop-all-gluster-processes.sh
@exec test -d ${PKG_SYSCONFDIR}/glusterd -a ! -d ${VARBASE}/lib/glusterd && ${MKDIR} -p ${VARBASE}/lib/ && mv ${PKG_SYSCONFDIR}/glusterd  ${VARBASE}/lib/glusterd && ${LN} -s ${VARBASE}/lib/glusterd  ${PKG_SYSCONFDIR}/glusterd
