# filesystems/glusterfs/PLIST
@comment $NetBSD: PLIST,v 1.32 2016/03/22 08:14:00 manu Exp $
${PYSITELIB}/gluster/__init__.py
${PYSITELIB}/gluster/__init__.pyc
${PYSITELIB}/gluster/__init__.pyo
${PYSITELIB}/gluster/glupy/__init__.py
${PYSITELIB}/gluster/glupy/__init__.pyc
${PYSITELIB}/gluster/glupy/__init__.pyo
bin/glusterfind
etc/ganesha/ganesha-ha.conf.sample
etc/glusterfs/gluster-rsyslog-5.8.conf
etc/glusterfs/gluster-rsyslog-7.2.conf
etc/glusterfs/glusterd.vol.sample
etc/glusterfs/glusterfs-georep-logrotate
etc/glusterfs/glusterfs-logrotate
etc/glusterfs/group-virt.example
etc/glusterfs/logger.conf.example
include/glusterfs/api/glfs-handles.h
include/glusterfs/api/glfs.h
include/glusterfs/gfchangelog/changelog.h
include/glusterfs/gfdb/gfdb_data_store.h
include/glusterfs/gfdb/gfdb_data_store_helper.h
include/glusterfs/gfdb/gfdb_data_store_types.h
include/glusterfs/gfdb/gfdb_mem-types.h
include/glusterfs/gfdb/gfdb_sqlite3.h
include/glusterfs/gfdb/gfdb_sqlite3_helper.h
lib/glusterfs/3.7.9/auth/addr.la
lib/glusterfs/3.7.9/auth/login.la
lib/glusterfs/3.7.9/rpc-transport/socket.la
lib/glusterfs/3.7.9/xlator/cluster/afr.la
lib/glusterfs/3.7.9/xlator/cluster/dht.la
lib/glusterfs/3.7.9/xlator/cluster/disperse.so
lib/glusterfs/3.7.9/xlator/cluster/distribute.so
lib/glusterfs/3.7.9/xlator/cluster/ec.la
lib/glusterfs/3.7.9/xlator/cluster/nufa.la
lib/glusterfs/3.7.9/xlator/cluster/pump.la
lib/glusterfs/3.7.9/xlator/cluster/replicate.so
lib/glusterfs/3.7.9/xlator/cluster/stripe.la
lib/glusterfs/3.7.9/xlator/cluster/switch.la
lib/glusterfs/3.7.9/xlator/cluster/tier.la
lib/glusterfs/3.7.9/xlator/debug/error-gen.la
lib/glusterfs/3.7.9/xlator/debug/io-stats.la
lib/glusterfs/3.7.9/xlator/debug/trace.la
lib/glusterfs/3.7.9/xlator/encryption/crypt.la
lib/glusterfs/3.7.9/xlator/encryption/rot-13.la
lib/glusterfs/3.7.9/xlator/features/access-control.so
lib/glusterfs/3.7.9/xlator/features/arbiter.la
lib/glusterfs/3.7.9/xlator/features/barrier.la
lib/glusterfs/3.7.9/xlator/features/bit-rot.la
lib/glusterfs/3.7.9/xlator/features/bitrot-stub.la
lib/glusterfs/3.7.9/xlator/features/cdc.la
lib/glusterfs/3.7.9/xlator/features/changelog.la
lib/glusterfs/3.7.9/xlator/features/changetimerecorder.la
lib/glusterfs/3.7.9/xlator/features/ganesha.la
lib/glusterfs/3.7.9/xlator/features/gfid-access.la
lib/glusterfs/3.7.9/xlator/features/glupy.la
lib/glusterfs/3.7.9/xlator/features/glupy/debug-trace.py
lib/glusterfs/3.7.9/xlator/features/glupy/debug-trace.pyc
lib/glusterfs/3.7.9/xlator/features/glupy/debug-trace.pyo
lib/glusterfs/3.7.9/xlator/features/glupy/helloworld.py
lib/glusterfs/3.7.9/xlator/features/glupy/helloworld.pyc
lib/glusterfs/3.7.9/xlator/features/glupy/helloworld.pyo
lib/glusterfs/3.7.9/xlator/features/glupy/negative.py
lib/glusterfs/3.7.9/xlator/features/glupy/negative.pyc
lib/glusterfs/3.7.9/xlator/features/glupy/negative.pyo
lib/glusterfs/3.7.9/xlator/features/index.la
lib/glusterfs/3.7.9/xlator/features/locks.la
lib/glusterfs/3.7.9/xlator/features/mac-compat.la
lib/glusterfs/3.7.9/xlator/features/marker.la
lib/glusterfs/3.7.9/xlator/features/posix-locks.so
lib/glusterfs/3.7.9/xlator/features/prot_client.la
lib/glusterfs/3.7.9/xlator/features/prot_dht.la
lib/glusterfs/3.7.9/xlator/features/prot_server.la
lib/glusterfs/3.7.9/xlator/features/quiesce.la
lib/glusterfs/3.7.9/xlator/features/quota.la
lib/glusterfs/3.7.9/xlator/features/quotad.la
lib/glusterfs/3.7.9/xlator/features/read-only.la
lib/glusterfs/3.7.9/xlator/features/shard.la
lib/glusterfs/3.7.9/xlator/features/snapview-client.la
lib/glusterfs/3.7.9/xlator/features/snapview-server.la
lib/glusterfs/3.7.9/xlator/features/trash.la
lib/glusterfs/3.7.9/xlator/features/upcall.la
lib/glusterfs/3.7.9/xlator/features/worm.la
lib/glusterfs/3.7.9/xlator/meta.la
lib/glusterfs/3.7.9/xlator/mgmt/glusterd.la
lib/glusterfs/3.7.9/xlator/mount/api.la
lib/glusterfs/3.7.9/xlator/mount/fuse.la
lib/glusterfs/3.7.9/xlator/nfs/server.la
lib/glusterfs/3.7.9/xlator/performance/io-cache.la
lib/glusterfs/3.7.9/xlator/performance/io-threads.la
lib/glusterfs/3.7.9/xlator/performance/md-cache.la
lib/glusterfs/3.7.9/xlator/performance/open-behind.la
lib/glusterfs/3.7.9/xlator/performance/quick-read.la
lib/glusterfs/3.7.9/xlator/performance/read-ahead.la
lib/glusterfs/3.7.9/xlator/performance/readdir-ahead.la
lib/glusterfs/3.7.9/xlator/performance/stat-prefetch.so
lib/glusterfs/3.7.9/xlator/performance/write-behind.la
lib/glusterfs/3.7.9/xlator/protocol/client.la
lib/glusterfs/3.7.9/xlator/protocol/server.la
lib/glusterfs/3.7.9/xlator/storage/posix.la
lib/glusterfs/3.7.9/xlator/system/posix-acl.la
lib/glusterfs/3.7.9/xlator/testing/features/template.la
lib/glusterfs/3.7.9/xlator/testing/performance/symlink-cache.la
lib/libgfapi.la
lib/libgfchangelog.la
lib/libgfdb.la
lib/libgfrpc.la
lib/libgfxdr.la
lib/libglusterfs.la
lib/ocf/resource.d/glusterfs/glusterd
lib/ocf/resource.d/glusterfs/volume
lib/ocf/resource.d/heartbeat/ganesha_grace
lib/ocf/resource.d/heartbeat/ganesha_mon
lib/ocf/resource.d/heartbeat/ganesha_nfsd
lib/pkgconfig/glusterfs-api.pc
lib/pkgconfig/libgfchangelog.pc
lib/pkgconfig/libgfdb.pc
libexec/ganesha/create-export-ganesha.sh
libexec/ganesha/dbus-send.sh
libexec/ganesha/ganesha-ha.sh
libexec/glusterfs/gfind_missing_files/gcrawler
libexec/glusterfs/gfind_missing_files/gfid_to_path.py
libexec/glusterfs/gfind_missing_files/gfid_to_path.sh
libexec/glusterfs/gfind_missing_files/gfind_missing_files.sh
libexec/glusterfs/glusterfind/__init__.py
libexec/glusterfs/glusterfind/__init__.pyc
libexec/glusterfs/glusterfind/__init__.pyo
libexec/glusterfs/glusterfind/S57glusterfind-delete-post.py
libexec/glusterfs/glusterfind/brickfind.py
libexec/glusterfs/glusterfind/changelog.py
libexec/glusterfs/glusterfind/changelogdata.py
libexec/glusterfs/glusterfind/changelogdata.pyc
libexec/glusterfs/glusterfind/changelogdata.pyo
libexec/glusterfs/glusterfind/conf.py
libexec/glusterfs/glusterfind/conf.pyc
libexec/glusterfs/glusterfind/conf.pyo
libexec/glusterfs/glusterfind/libgfchangelog.py
libexec/glusterfs/glusterfind/libgfchangelog.pyc
libexec/glusterfs/glusterfind/libgfchangelog.pyo
libexec/glusterfs/glusterfind/main.py
libexec/glusterfs/glusterfind/main.pyc
libexec/glusterfs/glusterfind/main.pyo
libexec/glusterfs/glusterfind/nodeagent.py
libexec/glusterfs/glusterfind/tool.conf
libexec/glusterfs/glusterfind/utils.py
libexec/glusterfs/glusterfind/utils.pyc
libexec/glusterfs/glusterfind/utils.pyo
libexec/glusterfs/gsyncd
libexec/glusterfs/gverify.sh
libexec/glusterfs/peer_add_secret_pub
libexec/glusterfs/peer_gsec_create
libexec/glusterfs/peer_mountbroker
libexec/glusterfs/python/syncdaemon/README.md
libexec/glusterfs/python/syncdaemon/__init__.py
libexec/glusterfs/python/syncdaemon/__init__.pyc
libexec/glusterfs/python/syncdaemon/__init__.pyo
libexec/glusterfs/python/syncdaemon/changelogagent.py
libexec/glusterfs/python/syncdaemon/changelogagent.pyc
libexec/glusterfs/python/syncdaemon/changelogagent.pyo
libexec/glusterfs/python/syncdaemon/configinterface.py
libexec/glusterfs/python/syncdaemon/configinterface.pyc
libexec/glusterfs/python/syncdaemon/configinterface.pyo
libexec/glusterfs/python/syncdaemon/gconf.py
libexec/glusterfs/python/syncdaemon/gconf.pyc
libexec/glusterfs/python/syncdaemon/gconf.pyo
libexec/glusterfs/python/syncdaemon/gsyncd.py
libexec/glusterfs/python/syncdaemon/gsyncd.pyc
libexec/glusterfs/python/syncdaemon/gsyncd.pyo
libexec/glusterfs/python/syncdaemon/gsyncdstatus.py
libexec/glusterfs/python/syncdaemon/gsyncdstatus.pyc
libexec/glusterfs/python/syncdaemon/gsyncdstatus.pyo
libexec/glusterfs/python/syncdaemon/ipaddr.py
libexec/glusterfs/python/syncdaemon/ipaddr.pyc
libexec/glusterfs/python/syncdaemon/ipaddr.pyo
libexec/glusterfs/python/syncdaemon/libcxattr.py
libexec/glusterfs/python/syncdaemon/libcxattr.pyc
libexec/glusterfs/python/syncdaemon/libcxattr.pyo
libexec/glusterfs/python/syncdaemon/libgfchangelog.py
libexec/glusterfs/python/syncdaemon/libgfchangelog.pyc
libexec/glusterfs/python/syncdaemon/libgfchangelog.pyo
libexec/glusterfs/python/syncdaemon/master.py
libexec/glusterfs/python/syncdaemon/master.pyc
libexec/glusterfs/python/syncdaemon/master.pyo
libexec/glusterfs/python/syncdaemon/monitor.py
libexec/glusterfs/python/syncdaemon/monitor.pyc
libexec/glusterfs/python/syncdaemon/monitor.pyo
libexec/glusterfs/python/syncdaemon/repce.py
libexec/glusterfs/python/syncdaemon/repce.pyc
libexec/glusterfs/python/syncdaemon/repce.pyo
libexec/glusterfs/python/syncdaemon/resource.py
libexec/glusterfs/python/syncdaemon/resource.pyc
libexec/glusterfs/python/syncdaemon/resource.pyo
libexec/glusterfs/python/syncdaemon/syncdutils.py
libexec/glusterfs/python/syncdaemon/syncdutils.pyc
libexec/glusterfs/python/syncdaemon/syncdutils.pyo
libexec/glusterfs/set_geo_rep_pem_keys.sh
man/man8/gluster.8
man/man8/glusterd.8
man/man8/glusterfs.8
man/man8/glusterfsd.8
man/man8/mount.glusterfs.8
sbin/gcron.py
sbin/gfind_missing_files
sbin/glfsheal
sbin/gluster
sbin/glusterd
sbin/glusterfs
sbin/glusterfsd
sbin/mount_glusterfs
sbin/snap_scheduler.py
sbin/umountd
share/doc/glusterfs/benchmarking/README
share/doc/glusterfs/benchmarking/glfs-bm.c
share/doc/glusterfs/benchmarking/launch-script.sh
share/doc/glusterfs/benchmarking/local-script.sh
share/doc/glusterfs/benchmarking/rdd.c
share/doc/glusterfs/glusterfs-mode.el
share/doc/glusterfs/glusterfs.vim
share/examples/rc.d/glusterd
share/glusterfs/scripts/generate-gfid-file.sh
share/glusterfs/scripts/get-gfid.sh
share/glusterfs/scripts/gsync-sync-gfid
share/glusterfs/scripts/gsync-upgrade.sh
share/glusterfs/scripts/post-upgrade-script-for-quota.sh
share/glusterfs/scripts/pre-upgrade-script-for-quota.sh
share/glusterfs/scripts/schedule_georep.py
share/glusterfs/scripts/slave-upgrade.sh
share/glusterfs/scripts/stop-all-gluster-processes.sh
@exec test -d ${PKG_SYSCONFDIR}/glusterd -a ! -d ${VARBASE}/lib/glusterd && ${MKDIR} -p ${VARBASE}/lib/ && mv ${PKG_SYSCONFDIR}/glusterd  ${VARBASE}/lib/glusterd && ${LN} -s ${VARBASE}/lib/glusterd  ${PKG_SYSCONFDIR}/glusterd
