# filesystems/glusterfs/PLIST
@comment $NetBSD: PLIST,v 1.21 2015/01/11 23:07:01 joerg Exp $
${PYSITELIB}/gluster/__init__.py
${PYSITELIB}/gluster/__init__.pyc
${PYSITELIB}/gluster/__init__.pyo
${PYSITELIB}/gluster/glupy.py
${PYSITELIB}/gluster/glupy.pyc
${PYSITELIB}/gluster/glupy.pyo
etc/glusterfs/gluster-rsyslog-5.8.conf
etc/glusterfs/gluster-rsyslog-7.2.conf
etc/glusterfs/glusterd.vol.sample
etc/glusterfs/glusterfs-georep-logrotate
etc/glusterfs/glusterfs-logrotate
etc/glusterfs/group-virt.example
etc/glusterfs/logger.conf.example
include/glusterfs/api/glfs-handles.h
include/glusterfs/api/glfs.h
include/glusterfs/gfchangelog/changelog.h
lib/glusterfs/3.6.0/auth/addr.la
lib/glusterfs/3.6.0/auth/login.la
lib/glusterfs/3.6.0/rpc-transport/socket.la
lib/glusterfs/3.6.0/xlator/cluster/afr.la
lib/glusterfs/3.6.0/xlator/cluster/dht.la
${PLIST.mmx}lib/glusterfs/3.6.0/xlator/cluster/disperse.so
lib/glusterfs/3.6.0/xlator/cluster/distribute.so
${PLIST.mmx}lib/glusterfs/3.6.0/xlator/cluster/ec.la
lib/glusterfs/3.6.0/xlator/cluster/nufa.la
lib/glusterfs/3.6.0/xlator/cluster/pump.la
lib/glusterfs/3.6.0/xlator/cluster/replicate.so
lib/glusterfs/3.6.0/xlator/cluster/stripe.la
lib/glusterfs/3.6.0/xlator/cluster/switch.la
lib/glusterfs/3.6.0/xlator/debug/error-gen.la
lib/glusterfs/3.6.0/xlator/debug/io-stats.la
lib/glusterfs/3.6.0/xlator/debug/trace.la
lib/glusterfs/3.6.0/xlator/encryption/crypt.la
lib/glusterfs/3.6.0/xlator/encryption/rot-13.la
lib/glusterfs/3.6.0/xlator/features/access-control.so
lib/glusterfs/3.6.0/xlator/features/barrier.la
lib/glusterfs/3.6.0/xlator/features/cdc.la
lib/glusterfs/3.6.0/xlator/features/changelog.la
lib/glusterfs/3.6.0/xlator/features/gfid-access.la
lib/glusterfs/3.6.0/xlator/features/glupy.la
lib/glusterfs/3.6.0/xlator/features/glupy/debug-trace.py
lib/glusterfs/3.6.0/xlator/features/glupy/debug-trace.pyc
lib/glusterfs/3.6.0/xlator/features/glupy/debug-trace.pyo
lib/glusterfs/3.6.0/xlator/features/glupy/helloworld.py
lib/glusterfs/3.6.0/xlator/features/glupy/helloworld.pyc
lib/glusterfs/3.6.0/xlator/features/glupy/helloworld.pyo
lib/glusterfs/3.6.0/xlator/features/glupy/negative.py
lib/glusterfs/3.6.0/xlator/features/glupy/negative.pyc
lib/glusterfs/3.6.0/xlator/features/glupy/negative.pyo
lib/glusterfs/3.6.0/xlator/features/index.la
lib/glusterfs/3.6.0/xlator/features/locks.la
lib/glusterfs/3.6.0/xlator/features/mac-compat.la
lib/glusterfs/3.6.0/xlator/features/marker.la
lib/glusterfs/3.6.0/xlator/features/posix-locks.so
lib/glusterfs/3.6.0/xlator/features/prot_client.la
lib/glusterfs/3.6.0/xlator/features/prot_dht.la
lib/glusterfs/3.6.0/xlator/features/prot_server.la
lib/glusterfs/3.6.0/xlator/features/quiesce.la
lib/glusterfs/3.6.0/xlator/features/quota.la
lib/glusterfs/3.6.0/xlator/features/quotad.la
lib/glusterfs/3.6.0/xlator/features/read-only.la
lib/glusterfs/3.6.0/xlator/features/snapview-client.la
lib/glusterfs/3.6.0/xlator/features/snapview-server.la
lib/glusterfs/3.6.0/xlator/features/worm.la
lib/glusterfs/3.6.0/xlator/meta.la
lib/glusterfs/3.6.0/xlator/mgmt/glusterd.la
lib/glusterfs/3.6.0/xlator/mount/api.la
lib/glusterfs/3.6.0/xlator/mount/fuse.la
lib/glusterfs/3.6.0/xlator/nfs/server.la
lib/glusterfs/3.6.0/xlator/performance/io-cache.la
lib/glusterfs/3.6.0/xlator/performance/io-threads.la
lib/glusterfs/3.6.0/xlator/performance/md-cache.la
lib/glusterfs/3.6.0/xlator/performance/open-behind.la
lib/glusterfs/3.6.0/xlator/performance/quick-read.la
lib/glusterfs/3.6.0/xlator/performance/read-ahead.la
lib/glusterfs/3.6.0/xlator/performance/readdir-ahead.la
lib/glusterfs/3.6.0/xlator/performance/stat-prefetch.so
lib/glusterfs/3.6.0/xlator/performance/write-behind.la
lib/glusterfs/3.6.0/xlator/protocol/client.la
lib/glusterfs/3.6.0/xlator/protocol/server.la
lib/glusterfs/3.6.0/xlator/storage/posix.la
lib/glusterfs/3.6.0/xlator/system/posix-acl.la
lib/glusterfs/3.6.0/xlator/testing/features/template.la
lib/glusterfs/3.6.0/xlator/testing/performance/symlink-cache.la
lib/libgfapi.la
lib/libgfchangelog.la
lib/libgfrpc.la
lib/libgfxdr.la
lib/libglusterfs.la
lib/ocf/resource.d/glusterfs/glusterd
lib/ocf/resource.d/glusterfs/volume
lib/pkgconfig/glusterfs-api.pc
lib/pkgconfig/libgfchangelog.pc
libexec/glusterfs/gsyncd
libexec/glusterfs/gverify.sh
libexec/glusterfs/peer_add_secret_pub
libexec/glusterfs/peer_gsec_create
libexec/glusterfs/python/syncdaemon/README.md
libexec/glusterfs/python/syncdaemon/__init__.py
libexec/glusterfs/python/syncdaemon/__init__.pyc
libexec/glusterfs/python/syncdaemon/__init__.pyo
libexec/glusterfs/python/syncdaemon/changelogagent.py
libexec/glusterfs/python/syncdaemon/changelogagent.pyc
libexec/glusterfs/python/syncdaemon/changelogagent.pyo
libexec/glusterfs/python/syncdaemon/configinterface.py
libexec/glusterfs/python/syncdaemon/configinterface.pyc
libexec/glusterfs/python/syncdaemon/configinterface.pyo
libexec/glusterfs/python/syncdaemon/gconf.py
libexec/glusterfs/python/syncdaemon/gconf.pyc
libexec/glusterfs/python/syncdaemon/gconf.pyo
libexec/glusterfs/python/syncdaemon/gsyncd.py
libexec/glusterfs/python/syncdaemon/gsyncd.pyc
libexec/glusterfs/python/syncdaemon/gsyncd.pyo
libexec/glusterfs/python/syncdaemon/ipaddr.py
libexec/glusterfs/python/syncdaemon/ipaddr.pyc
libexec/glusterfs/python/syncdaemon/ipaddr.pyo
libexec/glusterfs/python/syncdaemon/libcxattr.py
libexec/glusterfs/python/syncdaemon/libcxattr.pyc
libexec/glusterfs/python/syncdaemon/libcxattr.pyo
libexec/glusterfs/python/syncdaemon/libgfchangelog.py
libexec/glusterfs/python/syncdaemon/libgfchangelog.pyc
libexec/glusterfs/python/syncdaemon/libgfchangelog.pyo
libexec/glusterfs/python/syncdaemon/master.py
libexec/glusterfs/python/syncdaemon/master.pyc
libexec/glusterfs/python/syncdaemon/master.pyo
libexec/glusterfs/python/syncdaemon/monitor.py
libexec/glusterfs/python/syncdaemon/monitor.pyc
libexec/glusterfs/python/syncdaemon/monitor.pyo
libexec/glusterfs/python/syncdaemon/repce.py
libexec/glusterfs/python/syncdaemon/repce.pyc
libexec/glusterfs/python/syncdaemon/repce.pyo
libexec/glusterfs/python/syncdaemon/resource.py
libexec/glusterfs/python/syncdaemon/resource.pyc
libexec/glusterfs/python/syncdaemon/resource.pyo
libexec/glusterfs/python/syncdaemon/syncdutils.py
libexec/glusterfs/python/syncdaemon/syncdutils.pyc
libexec/glusterfs/python/syncdaemon/syncdutils.pyo
libexec/glusterfs/set_geo_rep_pem_keys.sh
man/man8/gluster.8
man/man8/glusterd.8
man/man8/glusterfs.8
man/man8/glusterfsd.8
man/man8/mount.glusterfs.8
sbin/gluster
sbin/glusterd
sbin/glusterfs
sbin/glusterfsd
sbin/mount_glusterfs
sbin/umountd
share/doc/glusterfs/benchmarking/README
share/doc/glusterfs/benchmarking/glfs-bm.c
share/doc/glusterfs/benchmarking/launch-script.sh
share/doc/glusterfs/benchmarking/local-script.sh
share/doc/glusterfs/benchmarking/rdd.c
share/doc/glusterfs/glusterfs-mode.el
share/doc/glusterfs/glusterfs.vim
share/examples/rc.d/glusterd
share/glusterfs/scripts/generate-gfid-file.sh
share/glusterfs/scripts/get-gfid.sh
share/glusterfs/scripts/gsync-sync-gfid
share/glusterfs/scripts/gsync-upgrade.sh
share/glusterfs/scripts/post-upgrade-script-for-quota.sh
share/glusterfs/scripts/pre-upgrade-script-for-quota.sh
share/glusterfs/scripts/slave-upgrade.sh
@exec test -d ${PKG_SYSCONFDIR}/glusterd -a ! -d ${VARBASE}/lib/glusterd && ${MKDIR} -p ${VARBASE}/lib/ && mv ${PKG_SYSCONFDIR}/glusterd  ${VARBASE}/lib/glusterd && ${LN} -s ${VARBASE}/lib/glusterd  ${PKG_SYSCONFDIR}/glusterd
