# chat/i2cbd/PLIST
@comment $NetBSD: PLIST,v 1.4 2006/05/21 16:08:32 rillig Exp $
man/man8/icbd.8
man/man8/icbdb.8
sbin/icbd
sbin/icbdb
share/examples/i2cbd/icbd.deny
share/examples/i2cbd/icbd_help
share/examples/i2cbd/motd
share/examples/i2cbd/news.1
share/examples/rc.d/icbd
@dirrm share/examples/i2cbd
