# net/ntp4/PLIST
@comment $NetBSD: PLIST,v 1.14 2014/01/12 17:01:02 spz Exp $
bin/sntp
man/man1/sntp.1
man/man5/ntp.conf.5
man/man5/ntp.keys.5
man/man8/ntp-keygen.8
man/man8/ntpd.8
man/man8/ntpdc.8
man/man8/ntpq.8
${PLIST.tickadj}sbin/calc_tickadj
sbin/ntp-keygen
sbin/ntp-wait
sbin/ntpd
sbin/ntpdate
sbin/ntpdc
sbin/ntpq
${PLIST.ntptime}sbin/ntptime
sbin/ntptrace
${PLIST.tickadj}sbin/tickadj
${PLIST.timetrim}sbin/timetrim
share/doc/ntp4/access.html
share/doc/ntp4/accopt.html
share/doc/ntp4/assoc.html
share/doc/ntp4/audio.html
share/doc/ntp4/authentic.html
share/doc/ntp4/authopt.html
share/doc/ntp4/autokey.html
share/doc/ntp4/bugs.html
share/doc/ntp4/build.html
share/doc/ntp4/clock.html
share/doc/ntp4/clockopt.html
share/doc/ntp4/cluster.html
share/doc/ntp4/comdex.html
share/doc/ntp4/config.html
share/doc/ntp4/confopt.html
share/doc/ntp4/copyright.html
share/doc/ntp4/debug.html
share/doc/ntp4/decode.html
share/doc/ntp4/discipline.html
share/doc/ntp4/discover.html
share/doc/ntp4/drivers/driver1.html
share/doc/ntp4/drivers/driver10.html
share/doc/ntp4/drivers/driver11.html
share/doc/ntp4/drivers/driver12.html
share/doc/ntp4/drivers/driver16.html
share/doc/ntp4/drivers/driver18.html
share/doc/ntp4/drivers/driver19.html
share/doc/ntp4/drivers/driver20.html
share/doc/ntp4/drivers/driver22.html
share/doc/ntp4/drivers/driver26.html
share/doc/ntp4/drivers/driver27.html
share/doc/ntp4/drivers/driver28.html
share/doc/ntp4/drivers/driver29.html
share/doc/ntp4/drivers/driver3.html
share/doc/ntp4/drivers/driver30.html
share/doc/ntp4/drivers/driver31.html
share/doc/ntp4/drivers/driver32.html
share/doc/ntp4/drivers/driver33.html
share/doc/ntp4/drivers/driver34.html
share/doc/ntp4/drivers/driver35.html
share/doc/ntp4/drivers/driver36.html
share/doc/ntp4/drivers/driver37.html
share/doc/ntp4/drivers/driver38.html
share/doc/ntp4/drivers/driver39.html
share/doc/ntp4/drivers/driver4.html
share/doc/ntp4/drivers/driver40.html
share/doc/ntp4/drivers/driver42.html
share/doc/ntp4/drivers/driver43.html
share/doc/ntp4/drivers/driver44.html
share/doc/ntp4/drivers/driver45.html
share/doc/ntp4/drivers/driver5.html
share/doc/ntp4/drivers/driver6.html
share/doc/ntp4/drivers/driver7.html
share/doc/ntp4/drivers/driver8.html
share/doc/ntp4/drivers/driver9.html
share/doc/ntp4/drivers/icons/home.gif
share/doc/ntp4/drivers/icons/mail2.gif
share/doc/ntp4/drivers/mx4200data.html
share/doc/ntp4/drivers/oncore-shmem.html
share/doc/ntp4/drivers/scripts/footer.txt
share/doc/ntp4/drivers/scripts/style.css
share/doc/ntp4/drivers/tf582_4.html
share/doc/ntp4/extern.html
share/doc/ntp4/filter.html
share/doc/ntp4/hints.html
share/doc/ntp4/hints/a-ux
share/doc/ntp4/hints/aix
share/doc/ntp4/hints/bsdi
share/doc/ntp4/hints/changes
share/doc/ntp4/hints/decosf1
share/doc/ntp4/hints/decosf2
share/doc/ntp4/hints/freebsd
share/doc/ntp4/hints/hpux
share/doc/ntp4/hints/linux
share/doc/ntp4/hints/mpeix
share/doc/ntp4/hints/notes-xntp-v3
share/doc/ntp4/hints/parse
share/doc/ntp4/hints/refclocks
share/doc/ntp4/hints/rs6000
share/doc/ntp4/hints/sco.html
share/doc/ntp4/hints/sgi
share/doc/ntp4/hints/solaris-dosynctodr.html
share/doc/ntp4/hints/solaris.html
share/doc/ntp4/hints/solaris.xtra.4023118
share/doc/ntp4/hints/solaris.xtra.4095849
share/doc/ntp4/hints/solaris.xtra.S99ntpd
share/doc/ntp4/hints/solaris.xtra.patchfreq
share/doc/ntp4/hints/sun4
share/doc/ntp4/hints/svr4-dell
share/doc/ntp4/hints/svr4_package
share/doc/ntp4/hints/todo
share/doc/ntp4/hints/vxworks.html
share/doc/ntp4/hints/winnt.html
share/doc/ntp4/history.html
share/doc/ntp4/howto.html
share/doc/ntp4/html/access.html
share/doc/ntp4/html/accopt.html
share/doc/ntp4/html/assoc.html
share/doc/ntp4/html/audio.html
share/doc/ntp4/html/authentic.html
share/doc/ntp4/html/authopt.html
share/doc/ntp4/html/autokey.html
share/doc/ntp4/html/bugs.html
share/doc/ntp4/html/build.html
share/doc/ntp4/html/clock.html
share/doc/ntp4/html/clockopt.html
share/doc/ntp4/html/cluster.html
share/doc/ntp4/html/comdex.html
share/doc/ntp4/html/config.html
share/doc/ntp4/html/confopt.html
share/doc/ntp4/html/copyright.html
share/doc/ntp4/html/debug.html
share/doc/ntp4/html/decode.html
share/doc/ntp4/html/discipline.html
share/doc/ntp4/html/discover.html
share/doc/ntp4/html/drivers/driver1.html
share/doc/ntp4/html/drivers/driver10.html
share/doc/ntp4/html/drivers/driver11.html
share/doc/ntp4/html/drivers/driver12.html
share/doc/ntp4/html/drivers/driver16.html
share/doc/ntp4/html/drivers/driver18.html
share/doc/ntp4/html/drivers/driver19.html
share/doc/ntp4/html/drivers/driver20.html
share/doc/ntp4/html/drivers/driver22.html
share/doc/ntp4/html/drivers/driver26.html
share/doc/ntp4/html/drivers/driver27.html
share/doc/ntp4/html/drivers/driver28.html
share/doc/ntp4/html/drivers/driver29.html
share/doc/ntp4/html/drivers/driver3.html
share/doc/ntp4/html/drivers/driver30.html
share/doc/ntp4/html/drivers/driver31.html
share/doc/ntp4/html/drivers/driver32.html
share/doc/ntp4/html/drivers/driver33.html
share/doc/ntp4/html/drivers/driver34.html
share/doc/ntp4/html/drivers/driver35.html
share/doc/ntp4/html/drivers/driver36.html
share/doc/ntp4/html/drivers/driver37.html
share/doc/ntp4/html/drivers/driver38.html
share/doc/ntp4/html/drivers/driver39.html
share/doc/ntp4/html/drivers/driver4.html
share/doc/ntp4/html/drivers/driver40.html
share/doc/ntp4/html/drivers/driver42.html
share/doc/ntp4/html/drivers/driver43.html
share/doc/ntp4/html/drivers/driver44.html
share/doc/ntp4/html/drivers/driver45.html
share/doc/ntp4/html/drivers/driver5.html
share/doc/ntp4/html/drivers/driver6.html
share/doc/ntp4/html/drivers/driver7.html
share/doc/ntp4/html/drivers/driver8.html
share/doc/ntp4/html/drivers/driver9.html
share/doc/ntp4/html/drivers/icons/home.gif
share/doc/ntp4/html/drivers/icons/mail2.gif
share/doc/ntp4/html/drivers/mx4200data.html
share/doc/ntp4/html/drivers/oncore-shmem.html
share/doc/ntp4/html/drivers/scripts/footer.txt
share/doc/ntp4/html/drivers/scripts/style.css
share/doc/ntp4/html/drivers/tf582_4.html
share/doc/ntp4/html/extern.html
share/doc/ntp4/html/filter.html
share/doc/ntp4/html/hints.html
share/doc/ntp4/html/hints/a-ux
share/doc/ntp4/html/hints/aix
share/doc/ntp4/html/hints/bsdi
share/doc/ntp4/html/hints/changes
share/doc/ntp4/html/hints/decosf1
share/doc/ntp4/html/hints/decosf2
share/doc/ntp4/html/hints/freebsd
share/doc/ntp4/html/hints/hpux
share/doc/ntp4/html/hints/linux
share/doc/ntp4/html/hints/mpeix
share/doc/ntp4/html/hints/notes-xntp-v3
share/doc/ntp4/html/hints/parse
share/doc/ntp4/html/hints/refclocks
share/doc/ntp4/html/hints/rs6000
share/doc/ntp4/html/hints/sco.html
share/doc/ntp4/html/hints/sgi
share/doc/ntp4/html/hints/solaris-dosynctodr.html
share/doc/ntp4/html/hints/solaris.html
share/doc/ntp4/html/hints/solaris.xtra.4023118
share/doc/ntp4/html/hints/solaris.xtra.4095849
share/doc/ntp4/html/hints/solaris.xtra.S99ntpd
share/doc/ntp4/html/hints/solaris.xtra.patchfreq
share/doc/ntp4/html/hints/sun4
share/doc/ntp4/html/hints/svr4-dell
share/doc/ntp4/html/hints/svr4_package
share/doc/ntp4/html/hints/todo
share/doc/ntp4/html/hints/vxworks.html
share/doc/ntp4/html/hints/winnt.html
share/doc/ntp4/html/history.html
share/doc/ntp4/html/howto.html
share/doc/ntp4/html/huffpuff.html
share/doc/ntp4/html/icons/home.gif
share/doc/ntp4/html/icons/mail2.gif
share/doc/ntp4/html/icons/sitemap.png
share/doc/ntp4/html/index.html
share/doc/ntp4/html/kern.html
share/doc/ntp4/html/kernpps.html
share/doc/ntp4/html/keygen.html
share/doc/ntp4/html/leap.html
share/doc/ntp4/html/miscopt.html
share/doc/ntp4/html/monopt.html
share/doc/ntp4/html/msyslog.html
share/doc/ntp4/html/ntp-wait.html
share/doc/ntp4/html/ntp_conf.html
share/doc/ntp4/html/ntpd.html
share/doc/ntp4/html/ntpdate.html
share/doc/ntp4/html/ntpdc.html
share/doc/ntp4/html/ntpdsim.html
share/doc/ntp4/html/ntpdsim_new.html
share/doc/ntp4/html/ntpq.html
share/doc/ntp4/html/ntptime.html
share/doc/ntp4/html/ntptrace.html
share/doc/ntp4/html/orphan.html
share/doc/ntp4/html/parsedata.html
share/doc/ntp4/html/parsenew.html
share/doc/ntp4/html/pic/9400n.jpg
share/doc/ntp4/html/pic/alice11.gif
share/doc/ntp4/html/pic/alice13.gif
share/doc/ntp4/html/pic/alice15.gif
share/doc/ntp4/html/pic/alice23.gif
share/doc/ntp4/html/pic/alice31.gif
share/doc/ntp4/html/pic/alice32.gif
share/doc/ntp4/html/pic/alice35.gif
share/doc/ntp4/html/pic/alice38.gif
share/doc/ntp4/html/pic/alice44.gif
share/doc/ntp4/html/pic/alice47.gif
share/doc/ntp4/html/pic/alice51.gif
share/doc/ntp4/html/pic/alice61.gif
share/doc/ntp4/html/pic/barnstable.gif
share/doc/ntp4/html/pic/beaver.gif
share/doc/ntp4/html/pic/boom3.gif
share/doc/ntp4/html/pic/boom3a.gif
share/doc/ntp4/html/pic/boom4.gif
share/doc/ntp4/html/pic/broad.gif
share/doc/ntp4/html/pic/bustardfly.gif
share/doc/ntp4/html/pic/c51.jpg
share/doc/ntp4/html/pic/description.jpg
share/doc/ntp4/html/pic/discipline.gif
share/doc/ntp4/html/pic/dogsnake.gif
share/doc/ntp4/html/pic/driver29.gif
share/doc/ntp4/html/pic/driver43_1.gif
share/doc/ntp4/html/pic/driver43_2.jpg
share/doc/ntp4/html/pic/fg6021.gif
share/doc/ntp4/html/pic/fg6039.jpg
share/doc/ntp4/html/pic/fig_3_1.gif
share/doc/ntp4/html/pic/flatheads.gif
share/doc/ntp4/html/pic/flt1.gif
share/doc/ntp4/html/pic/flt2.gif
share/doc/ntp4/html/pic/flt3.gif
share/doc/ntp4/html/pic/flt4.gif
share/doc/ntp4/html/pic/flt5.gif
share/doc/ntp4/html/pic/flt6.gif
share/doc/ntp4/html/pic/flt7.gif
share/doc/ntp4/html/pic/flt8.gif
share/doc/ntp4/html/pic/flt9.gif
share/doc/ntp4/html/pic/freq1211.gif
share/doc/ntp4/html/pic/gadget.jpg
share/doc/ntp4/html/pic/gps167.jpg
share/doc/ntp4/html/pic/group.gif
share/doc/ntp4/html/pic/hornraba.gif
share/doc/ntp4/html/pic/igclock.gif
share/doc/ntp4/html/pic/neoclock4x.gif
share/doc/ntp4/html/pic/offset1211.gif
share/doc/ntp4/html/pic/oncore_evalbig.gif
share/doc/ntp4/html/pic/oncore_remoteant.jpg
share/doc/ntp4/html/pic/oncore_utplusbig.gif
share/doc/ntp4/html/pic/oz2.gif
share/doc/ntp4/html/pic/panda.gif
share/doc/ntp4/html/pic/pd_om006.gif
share/doc/ntp4/html/pic/pd_om011.gif
share/doc/ntp4/html/pic/peer.gif
share/doc/ntp4/html/pic/pogo.gif
share/doc/ntp4/html/pic/pogo1a.gif
share/doc/ntp4/html/pic/pogo3a.gif
share/doc/ntp4/html/pic/pogo4.gif
share/doc/ntp4/html/pic/pogo5.gif
share/doc/ntp4/html/pic/pogo6.gif
share/doc/ntp4/html/pic/pogo7.gif
share/doc/ntp4/html/pic/pogo8.gif
share/doc/ntp4/html/pic/pzf509.jpg
share/doc/ntp4/html/pic/pzf511.jpg
share/doc/ntp4/html/pic/rabbit.gif
share/doc/ntp4/html/pic/radio2.jpg
share/doc/ntp4/html/pic/sheepb.jpg
share/doc/ntp4/html/pic/stack1a.jpg
share/doc/ntp4/html/pic/sx5.gif
share/doc/ntp4/html/pic/thunderbolt.jpg
share/doc/ntp4/html/pic/tonea.gif
share/doc/ntp4/html/pic/tribeb.gif
share/doc/ntp4/html/pic/wingdorothy.gif
share/doc/ntp4/html/poll.html
share/doc/ntp4/html/pps.html
share/doc/ntp4/html/prefer.html
share/doc/ntp4/html/quick.html
share/doc/ntp4/html/rate.html
share/doc/ntp4/html/rdebug.html
share/doc/ntp4/html/refclock.html
share/doc/ntp4/html/release.html
share/doc/ntp4/html/scripts/accopt.txt
share/doc/ntp4/html/scripts/audio.txt
share/doc/ntp4/html/scripts/authopt.txt
share/doc/ntp4/html/scripts/clockopt.txt
share/doc/ntp4/html/scripts/command.txt
share/doc/ntp4/html/scripts/config.txt
share/doc/ntp4/html/scripts/confopt.txt
share/doc/ntp4/html/scripts/external.txt
share/doc/ntp4/html/scripts/footer.txt
share/doc/ntp4/html/scripts/hand.txt
share/doc/ntp4/html/scripts/install.txt
share/doc/ntp4/html/scripts/manual.txt
share/doc/ntp4/html/scripts/misc.txt
share/doc/ntp4/html/scripts/miscopt.txt
share/doc/ntp4/html/scripts/monopt.txt
share/doc/ntp4/html/scripts/refclock.txt
share/doc/ntp4/html/scripts/special.txt
share/doc/ntp4/html/scripts/style.css
share/doc/ntp4/html/select.html
share/doc/ntp4/html/sitemap.html
share/doc/ntp4/html/sntp.html
share/doc/ntp4/html/stats.html
share/doc/ntp4/html/tickadj.html
share/doc/ntp4/html/warp.html
share/doc/ntp4/html/xleave.html
share/doc/ntp4/huffpuff.html
share/doc/ntp4/icons/home.gif
share/doc/ntp4/icons/mail2.gif
share/doc/ntp4/icons/sitemap.png
share/doc/ntp4/index.html
share/doc/ntp4/kern.html
share/doc/ntp4/kernpps.html
share/doc/ntp4/keygen.html
share/doc/ntp4/leap.html
share/doc/ntp4/miscopt.html
share/doc/ntp4/monopt.html
share/doc/ntp4/msyslog.html
share/doc/ntp4/ntp-keygen.html
share/doc/ntp4/ntp-wait.html
share/doc/ntp4/ntp.conf.html
share/doc/ntp4/ntp.keys.html
share/doc/ntp4/ntp_conf.html
share/doc/ntp4/ntpd.html
share/doc/ntp4/ntpdate.html
share/doc/ntp4/ntpdc.html
share/doc/ntp4/ntpdsim.html
share/doc/ntp4/ntpdsim_new.html
share/doc/ntp4/ntpq.html
share/doc/ntp4/ntpsnmpd.html
share/doc/ntp4/ntpsweep.html
share/doc/ntp4/ntptime.html
share/doc/ntp4/ntptrace.html
share/doc/ntp4/orphan.html
share/doc/ntp4/parsedata.html
share/doc/ntp4/parsenew.html
share/doc/ntp4/pic/9400n.jpg
share/doc/ntp4/pic/alice11.gif
share/doc/ntp4/pic/alice13.gif
share/doc/ntp4/pic/alice15.gif
share/doc/ntp4/pic/alice23.gif
share/doc/ntp4/pic/alice31.gif
share/doc/ntp4/pic/alice32.gif
share/doc/ntp4/pic/alice35.gif
share/doc/ntp4/pic/alice38.gif
share/doc/ntp4/pic/alice44.gif
share/doc/ntp4/pic/alice47.gif
share/doc/ntp4/pic/alice51.gif
share/doc/ntp4/pic/alice61.gif
share/doc/ntp4/pic/barnstable.gif
share/doc/ntp4/pic/beaver.gif
share/doc/ntp4/pic/boom3.gif
share/doc/ntp4/pic/boom3a.gif
share/doc/ntp4/pic/boom4.gif
share/doc/ntp4/pic/broad.gif
share/doc/ntp4/pic/bustardfly.gif
share/doc/ntp4/pic/c51.jpg
share/doc/ntp4/pic/description.jpg
share/doc/ntp4/pic/discipline.gif
share/doc/ntp4/pic/dogsnake.gif
share/doc/ntp4/pic/driver29.gif
share/doc/ntp4/pic/driver43_1.gif
share/doc/ntp4/pic/driver43_2.jpg
share/doc/ntp4/pic/fg6021.gif
share/doc/ntp4/pic/fg6039.jpg
share/doc/ntp4/pic/fig_3_1.gif
share/doc/ntp4/pic/flatheads.gif
share/doc/ntp4/pic/flt1.gif
share/doc/ntp4/pic/flt2.gif
share/doc/ntp4/pic/flt3.gif
share/doc/ntp4/pic/flt4.gif
share/doc/ntp4/pic/flt5.gif
share/doc/ntp4/pic/flt6.gif
share/doc/ntp4/pic/flt7.gif
share/doc/ntp4/pic/flt8.gif
share/doc/ntp4/pic/flt9.gif
share/doc/ntp4/pic/freq1211.gif
share/doc/ntp4/pic/gadget.jpg
share/doc/ntp4/pic/gps167.jpg
share/doc/ntp4/pic/group.gif
share/doc/ntp4/pic/hornraba.gif
share/doc/ntp4/pic/igclock.gif
share/doc/ntp4/pic/neoclock4x.gif
share/doc/ntp4/pic/offset1211.gif
share/doc/ntp4/pic/oncore_evalbig.gif
share/doc/ntp4/pic/oncore_remoteant.jpg
share/doc/ntp4/pic/oncore_utplusbig.gif
share/doc/ntp4/pic/oz2.gif
share/doc/ntp4/pic/panda.gif
share/doc/ntp4/pic/pd_om006.gif
share/doc/ntp4/pic/pd_om011.gif
share/doc/ntp4/pic/peer.gif
share/doc/ntp4/pic/pogo.gif
share/doc/ntp4/pic/pogo1a.gif
share/doc/ntp4/pic/pogo3a.gif
share/doc/ntp4/pic/pogo4.gif
share/doc/ntp4/pic/pogo5.gif
share/doc/ntp4/pic/pogo6.gif
share/doc/ntp4/pic/pogo7.gif
share/doc/ntp4/pic/pogo8.gif
share/doc/ntp4/pic/pzf509.jpg
share/doc/ntp4/pic/pzf511.jpg
share/doc/ntp4/pic/rabbit.gif
share/doc/ntp4/pic/radio2.jpg
share/doc/ntp4/pic/sheepb.jpg
share/doc/ntp4/pic/stack1a.jpg
share/doc/ntp4/pic/sx5.gif
share/doc/ntp4/pic/thunderbolt.jpg
share/doc/ntp4/pic/tonea.gif
share/doc/ntp4/pic/tribeb.gif
share/doc/ntp4/pic/wingdorothy.gif
share/doc/ntp4/poll.html
share/doc/ntp4/pps.html
share/doc/ntp4/prefer.html
share/doc/ntp4/quick.html
share/doc/ntp4/rate.html
share/doc/ntp4/rdebug.html
share/doc/ntp4/refclock.html
share/doc/ntp4/release.html
share/doc/ntp4/scripts/accopt.txt
share/doc/ntp4/scripts/audio.txt
share/doc/ntp4/scripts/authopt.txt
share/doc/ntp4/scripts/clockopt.txt
share/doc/ntp4/scripts/command.txt
share/doc/ntp4/scripts/config.txt
share/doc/ntp4/scripts/confopt.txt
share/doc/ntp4/scripts/external.txt
share/doc/ntp4/scripts/footer.txt
share/doc/ntp4/scripts/hand.txt
share/doc/ntp4/scripts/install.txt
share/doc/ntp4/scripts/manual.txt
share/doc/ntp4/scripts/misc.txt
share/doc/ntp4/scripts/miscopt.txt
share/doc/ntp4/scripts/monopt.txt
share/doc/ntp4/scripts/refclock.txt
share/doc/ntp4/scripts/special.txt
share/doc/ntp4/scripts/style.css
share/doc/ntp4/select.html
share/doc/ntp4/sitemap.html
share/doc/ntp4/sntp.html
share/doc/ntp4/stats.html
share/doc/ntp4/tickadj.html
share/doc/ntp4/warp.html
share/doc/ntp4/xleave.html
share/doc/sntp/sntp.html
share/examples/ntp4/README
share/examples/ntp4/baldwin.conf
share/examples/ntp4/beauregard.conf
share/examples/ntp4/grundoon.conf
share/examples/ntp4/malarky.conf
share/examples/ntp4/pogo.conf
share/examples/ntp4/rackety.conf
share/examples/rc.d/ntpd
share/examples/rc.d/ntpdate
share/ntp/lib/NTP/Util.pm
@pkgdir include/event2
