# net/ntp4/PLIST
@comment $NetBSD: PLIST,v 1.13 2009/06/14 18:09:37 joerg Exp $
bin/ntp-keygen
bin/ntp-wait
bin/ntpd
bin/ntpdate
bin/ntpdc
bin/ntpq
${PLIST.ntptime}bin/ntptime
bin/ntptrace
bin/sntp
${PLIST.tickadj}bin/tickadj
${PLIST.timetrim}bin/timetrim
man/man1/ntp-keygen.1
man/man1/ntpd.1
man/man1/ntpdc.1
man/man1/ntpdsim.1
man/man1/ntpq.1
man/man1/sntp.1
share/doc/ntp4/accopt.html
share/doc/ntp4/assoc.html
share/doc/ntp4/audio.html
share/doc/ntp4/authopt.html
share/doc/ntp4/build/build.html
share/doc/ntp4/build/config.html
share/doc/ntp4/build/hints.html
share/doc/ntp4/build/hints/a-ux
share/doc/ntp4/build/hints/aix
share/doc/ntp4/build/hints/bsdi
share/doc/ntp4/build/hints/changes
share/doc/ntp4/build/hints/decosf1
share/doc/ntp4/build/hints/decosf2
share/doc/ntp4/build/hints/freebsd
share/doc/ntp4/build/hints/hpux
share/doc/ntp4/build/hints/linux
share/doc/ntp4/build/hints/mpeix
share/doc/ntp4/build/hints/netbsd
share/doc/ntp4/build/hints/notes-xntp-v3
share/doc/ntp4/build/hints/parse
share/doc/ntp4/build/hints/refclocks
share/doc/ntp4/build/hints/rs6000
share/doc/ntp4/build/hints/sco.html
share/doc/ntp4/build/hints/sgi
share/doc/ntp4/build/hints/solaris-dosynctodr.html
share/doc/ntp4/build/hints/solaris.html
share/doc/ntp4/build/hints/solaris.xtra.4023118
share/doc/ntp4/build/hints/solaris.xtra.4095849
share/doc/ntp4/build/hints/solaris.xtra.S99ntpd
share/doc/ntp4/build/hints/solaris.xtra.patchfreq
share/doc/ntp4/build/hints/sun4
share/doc/ntp4/build/hints/svr4-dell
share/doc/ntp4/build/hints/svr4_package
share/doc/ntp4/build/hints/todo
share/doc/ntp4/build/hints/vxworks.html
share/doc/ntp4/build/hints/winnt.html
share/doc/ntp4/build/patches.html
share/doc/ntp4/build/porting.html
share/doc/ntp4/build/quick.html
share/doc/ntp4/build/scripts/footer.txt
share/doc/ntp4/build/scripts/links10.txt
share/doc/ntp4/build/scripts/links11.txt
share/doc/ntp4/build/scripts/links12.txt
share/doc/ntp4/build/scripts/links7.txt
share/doc/ntp4/build/scripts/links8.txt
share/doc/ntp4/build/scripts/links9.txt
share/doc/ntp4/build/scripts/style.css
share/doc/ntp4/clockopt.html
share/doc/ntp4/confopt.html
share/doc/ntp4/copyright.html
share/doc/ntp4/debug.html
share/doc/ntp4/drivers/driver1.html
share/doc/ntp4/drivers/driver10.html
share/doc/ntp4/drivers/driver11.html
share/doc/ntp4/drivers/driver12.html
share/doc/ntp4/drivers/driver16.html
share/doc/ntp4/drivers/driver18.html
share/doc/ntp4/drivers/driver19.html
share/doc/ntp4/drivers/driver2.html
share/doc/ntp4/drivers/driver20.html
share/doc/ntp4/drivers/driver22.html
share/doc/ntp4/drivers/driver26.html
share/doc/ntp4/drivers/driver27.html
share/doc/ntp4/drivers/driver28.html
share/doc/ntp4/drivers/driver29.html
share/doc/ntp4/drivers/driver3.html
share/doc/ntp4/drivers/driver30.html
share/doc/ntp4/drivers/driver31.html
share/doc/ntp4/drivers/driver32.html
share/doc/ntp4/drivers/driver33.html
share/doc/ntp4/drivers/driver34.html
share/doc/ntp4/drivers/driver35.html
share/doc/ntp4/drivers/driver36.html
share/doc/ntp4/drivers/driver37.html
share/doc/ntp4/drivers/driver38.html
share/doc/ntp4/drivers/driver39.html
share/doc/ntp4/drivers/driver4.html
share/doc/ntp4/drivers/driver40.html
share/doc/ntp4/drivers/driver42.html
share/doc/ntp4/drivers/driver43.html
share/doc/ntp4/drivers/driver44.html
share/doc/ntp4/drivers/driver5.html
share/doc/ntp4/drivers/driver6.html
share/doc/ntp4/drivers/driver7.html
share/doc/ntp4/drivers/driver8.html
share/doc/ntp4/drivers/driver9.html
share/doc/ntp4/drivers/icons/home.gif
share/doc/ntp4/drivers/icons/mail2.gif
share/doc/ntp4/drivers/oncore-shmem.html
share/doc/ntp4/drivers/scripts/footer.txt
share/doc/ntp4/drivers/scripts/style.css
share/doc/ntp4/drivers/tf582_4.html
share/doc/ntp4/extern.html
share/doc/ntp4/gadget.html
share/doc/ntp4/groups.html
share/doc/ntp4/howto.html
share/doc/ntp4/icons/home.gif
share/doc/ntp4/icons/mail2.gif
share/doc/ntp4/index.html
share/doc/ntp4/kern.html
share/doc/ntp4/keygen.html
share/doc/ntp4/ldisc.html
share/doc/ntp4/manyopt.html
share/doc/ntp4/measure.html
share/doc/ntp4/miscopt.html
share/doc/ntp4/monopt.html
share/doc/ntp4/msyslog.html
share/doc/ntp4/mx4200data.html
share/doc/ntp4/notes.html
share/doc/ntp4/ntp_conf.html
share/doc/ntp4/ntpd.html
share/doc/ntp4/ntpdate.html
share/doc/ntp4/ntpdc.html
share/doc/ntp4/ntpdsim.html
share/doc/ntp4/ntpdsim_new.html
share/doc/ntp4/ntpq.html
share/doc/ntp4/ntptime.html
share/doc/ntp4/ntptrace.html
share/doc/ntp4/parsedata.html
share/doc/ntp4/parsenew.html
share/doc/ntp4/pic/9400n.jpg
share/doc/ntp4/pic/alice11.gif
share/doc/ntp4/pic/alice13.gif
share/doc/ntp4/pic/alice23.gif
share/doc/ntp4/pic/alice31.gif
share/doc/ntp4/pic/alice32.gif
share/doc/ntp4/pic/alice35.gif
share/doc/ntp4/pic/alice38.gif
share/doc/ntp4/pic/alice44.gif
share/doc/ntp4/pic/alice47.gif
share/doc/ntp4/pic/alice51.gif
share/doc/ntp4/pic/alice61.gif
share/doc/ntp4/pic/barnstable.gif
share/doc/ntp4/pic/beaver.gif
share/doc/ntp4/pic/boom3.gif
share/doc/ntp4/pic/boom3a.gif
share/doc/ntp4/pic/bustardfly.gif
share/doc/ntp4/pic/c51.jpg
share/doc/ntp4/pic/dogsnake.gif
share/doc/ntp4/pic/driver29.gif
share/doc/ntp4/pic/driver43_1.gif
share/doc/ntp4/pic/driver43_2.jpg
share/doc/ntp4/pic/fg6021.gif
share/doc/ntp4/pic/fg6039.jpg
share/doc/ntp4/pic/gadget.jpg
share/doc/ntp4/pic/gps167.jpg
share/doc/ntp4/pic/hornraba.gif
share/doc/ntp4/pic/igclock.gif
share/doc/ntp4/pic/neoclock4x.gif
share/doc/ntp4/pic/oncore_evalbig.gif
share/doc/ntp4/pic/oncore_remoteant.jpg
share/doc/ntp4/pic/oncore_utplusbig.gif
share/doc/ntp4/pic/oz2.gif
share/doc/ntp4/pic/panda.gif
share/doc/ntp4/pic/pd_om006.gif
share/doc/ntp4/pic/pd_om011.gif
share/doc/ntp4/pic/pogo.gif
share/doc/ntp4/pic/pogo1a.gif
share/doc/ntp4/pic/pogo3a.gif
share/doc/ntp4/pic/pogo4.gif
share/doc/ntp4/pic/pogo5.gif
share/doc/ntp4/pic/pogo6.gif
share/doc/ntp4/pic/pogo8.gif
share/doc/ntp4/pic/pzf511.jpg
share/doc/ntp4/pic/rabbit.gif
share/doc/ntp4/pic/radio2.jpg
share/doc/ntp4/pic/sheepb.jpg
share/doc/ntp4/pic/stack1a.jpg
share/doc/ntp4/pic/tonea.gif
share/doc/ntp4/pic/wingdorothy.gif
share/doc/ntp4/pps.html
share/doc/ntp4/prefer.html
share/doc/ntp4/rdebug.html
share/doc/ntp4/refclock.html
share/doc/ntp4/release.html
share/doc/ntp4/scripts/footer.txt
share/doc/ntp4/scripts/links10.txt
share/doc/ntp4/scripts/links11.txt
share/doc/ntp4/scripts/links12.txt
share/doc/ntp4/scripts/links7.txt
share/doc/ntp4/scripts/links8.txt
share/doc/ntp4/scripts/links9.txt
share/doc/ntp4/scripts/style.css
share/doc/ntp4/sntp.html
share/doc/ntp4/tickadj.html
share/examples/ntp4/README
share/examples/ntp4/baldwin.conf
share/examples/ntp4/beauregard.conf
share/examples/ntp4/grundoon.conf
share/examples/ntp4/malarky.conf
share/examples/ntp4/pogo.conf
share/examples/ntp4/rackety.conf
share/examples/rc.d/ntpd
share/examples/rc.d/ntpdate
