# pkgtools/tinderbox-dragonfly/PLIST
@comment $NetBSD: PLIST,v 1.1 2011/11/02 21:13:16 marino Exp $
man/man1/tc-configCcache.1
man/man1/tc-configDistfile.1
man/man1/tc-configGet.1
man/man1/tc-configJail.1
man/man1/tc-configTinderd.1
man/man1/tc-init.1
share/examples/rc.d/tinderd
tinderbox/scripts/.version
tinderbox/scripts/README
tinderbox/scripts/ds.ph.dist
tinderbox/scripts/contrib/hooks/make-jobs-safe/README.txt
tinderbox/scripts/contrib/hooks/make-jobs-safe/make-jobs-safe.sh
tinderbox/scripts/contrib/hooks/kse-fix/libmap.conf
tinderbox/scripts/contrib/hooks/kse-fix/README.txt
tinderbox/scripts/contrib/hooks/kse-fix/kse-fix.sh
tinderbox/scripts/etc/env/.keep_me
tinderbox/scripts/lib/Tinderbox/Build.pm
tinderbox/scripts/lib/Tinderbox/BuildPortsQueue.pm
tinderbox/scripts/lib/Tinderbox/Config.pm
tinderbox/scripts/lib/Tinderbox/Hook.pm
tinderbox/scripts/lib/Tinderbox/Jail.pm
tinderbox/scripts/lib/Tinderbox/MakeCache.pm
tinderbox/scripts/lib/Tinderbox/Port.pm
tinderbox/scripts/lib/Tinderbox/PortFailPattern.pm
tinderbox/scripts/lib/Tinderbox/PortFailReason.pm
tinderbox/scripts/lib/Tinderbox/PortsTree.pm
tinderbox/scripts/lib/Tinderbox/TinderObject.pm
tinderbox/scripts/lib/Tinderbox/TinderboxDS.pm
tinderbox/scripts/lib/Tinderbox/User.pm
tinderbox/scripts/lib/buildscript
tinderbox/scripts/lib/db-mysql.sh
tinderbox/scripts/lib/db-pgsql.sh
tinderbox/scripts/lib/enterbuild
tinderbox/scripts/lib/makemake
tinderbox/scripts/lib/pnohang.c
tinderbox/scripts/lib/portbuild
tinderbox/scripts/lib/tc_command.pl
tinderbox/scripts/lib/tc_command.sh
tinderbox/scripts/lib/tinderbox.env
tinderbox/scripts/lib/tinderlib.pl
tinderbox/scripts/lib/tinderlib.sh
tinderbox/scripts/sql/genschema
tinderbox/scripts/sql/schema.mysql.post
tinderbox/scripts/sql/schema.mysql.pre
tinderbox/scripts/sql/schema.pgsql.post
tinderbox/scripts/sql/schema.pgsql.pre
tinderbox/scripts/sql/values.config
tinderbox/scripts/sql/values.hooks
tinderbox/scripts/sql/values.lp
tinderbox/scripts/sql/values.pfp
tinderbox/scripts/sql/values.pfr
tinderbox/scripts/tc
tinderbox/scripts/tinderbox.ph.dist
tinderbox/scripts/tinderd
tinderbox/scripts/upgrade/build_ports.map
tinderbox/scripts/upgrade/build_ports_queue.map
tinderbox/scripts/upgrade/builds.map
tinderbox/scripts/upgrade/config.map
tinderbox/scripts/upgrade/hooks.map
tinderbox/scripts/upgrade/jails.map
tinderbox/scripts/upgrade/mig_mysql_tinderbox-3.0_to_3.1.sql
tinderbox/scripts/upgrade/mig_mysql_tinderbox-3.1_to_3.2.sql
tinderbox/scripts/upgrade/mig_pgsql_tinderbox-3.0_to_3.1.sql
tinderbox/scripts/upgrade/mig_pgsql_tinderbox-3.1_to_3.2.sql
tinderbox/scripts/upgrade/mig_pgsql_tinderbox-3.2_to_3.2.1.sql
tinderbox/scripts/upgrade/mig_mysql_tinderbox-3.2_to_3.2.1.sql
tinderbox/scripts/upgrade/order.lst
tinderbox/scripts/upgrade/user_permissions.map
${PLIST.WEBUI}tinderbox/scripts/webui/core/Build.php
${PLIST.WEBUI}tinderbox/scripts/webui/core/BuildGroups.php
${PLIST.WEBUI}tinderbox/scripts/webui/core/BuildPortsQueue.php
${PLIST.WEBUI}tinderbox/scripts/webui/core/Config.php
${PLIST.WEBUI}tinderbox/scripts/webui/core/Hooks.php
${PLIST.WEBUI}tinderbox/scripts/webui/core/Jail.php
${PLIST.WEBUI}tinderbox/scripts/webui/core/LogfilePattern.php
${PLIST.WEBUI}tinderbox/scripts/webui/core/Port.php
${PLIST.WEBUI}tinderbox/scripts/webui/core/PortFailPattern.php
${PLIST.WEBUI}tinderbox/scripts/webui/core/PortFailReason.php
${PLIST.WEBUI}tinderbox/scripts/webui/core/PortsTree.php
${PLIST.WEBUI}tinderbox/scripts/webui/core/TinderObject.php
${PLIST.WEBUI}tinderbox/scripts/webui/core/TinderboxDS.php
${PLIST.WEBUI}tinderbox/scripts/webui/core/User.php
${PLIST.WEBUI}tinderbox/scripts/webui/core/functions.php
${PLIST.WEBUI}tinderbox/scripts/webui/inc_ds.php.dist
${PLIST.WEBUI}tinderbox/scripts/webui/inc_tinderbox.php.dist
${PLIST.WEBUI}tinderbox/scripts/webui/index.php
${PLIST.WEBUI}tinderbox/scripts/webui/favicon.ico
${PLIST.WEBUI}tinderbox/scripts/webui/module/module.php
${PLIST.WEBUI}tinderbox/scripts/webui/module/moduleBuildGroups.php
${PLIST.WEBUI}tinderbox/scripts/webui/module/moduleBuildPorts.php
${PLIST.WEBUI}tinderbox/scripts/webui/module/moduleBuilds.php
${PLIST.WEBUI}tinderbox/scripts/webui/module/moduleConfig.php
${PLIST.WEBUI}tinderbox/scripts/webui/module/moduleLogs.php
${PLIST.WEBUI}tinderbox/scripts/webui/module/modulePortFailureReasons.php
${PLIST.WEBUI}tinderbox/scripts/webui/module/modulePorts.php
${PLIST.WEBUI}tinderbox/scripts/webui/module/moduleRss.php
${PLIST.WEBUI}tinderbox/scripts/webui/module/moduleSession.php
${PLIST.WEBUI}tinderbox/scripts/webui/module/moduleTinderd.php
${PLIST.WEBUI}tinderbox/scripts/webui/module/moduleUsers.php
${PLIST.WEBUI}tinderbox/scripts/webui/templates/default/config.tpl
${PLIST.WEBUI}tinderbox/scripts/webui/templates/default/current_buildports.tpl
${PLIST.WEBUI}tinderbox/scripts/webui/templates/default/describe_port.tpl
${PLIST.WEBUI}tinderbox/scripts/webui/templates/default/display_login.tpl
${PLIST.WEBUI}tinderbox/scripts/webui/templates/default/display_build_groups.tpl
${PLIST.WEBUI}tinderbox/scripts/webui/templates/default/display_markup_log.tpl
${PLIST.WEBUI}tinderbox/scripts/webui/templates/default/failed_buildports.tpl
${PLIST.WEBUI}tinderbox/scripts/webui/templates/default/latest_buildports.tpl
${PLIST.WEBUI}tinderbox/scripts/webui/templates/default/list_buildports.tpl
${PLIST.WEBUI}tinderbox/scripts/webui/templates/default/list_builds.tpl
${PLIST.WEBUI}tinderbox/scripts/webui/templates/default/list_failure_reasons.tpl
${PLIST.WEBUI}tinderbox/scripts/webui/templates/default/list_tinderd_queue.tpl
${PLIST.WEBUI}tinderbox/scripts/webui/templates/default/messages.inc
${PLIST.WEBUI}tinderbox/scripts/webui/templates/default/please_login.tpl
${PLIST.WEBUI}tinderbox/scripts/webui/templates/default/rss.tpl
${PLIST.WEBUI}tinderbox/scripts/webui/templates/default/tinderstyle.css
${PLIST.WEBUI}tinderbox/scripts/webui/templates/default/tinderstyle.js
${PLIST.WEBUI}tinderbox/scripts/webui/templates/default/user_admin.tpl
${PLIST.WEBUI}tinderbox/scripts/webui/templates/default/user_permissions.tpl
${PLIST.WEBUI}tinderbox/scripts/webui/templates/default/user_properties.tpl
${PLIST.WEBUI}tinderbox/scripts/webui/templates/paefchen/config.tpl
${PLIST.WEBUI}tinderbox/scripts/webui/templates/paefchen/current_buildports.tpl
${PLIST.WEBUI}tinderbox/scripts/webui/templates/paefchen/describe_port.tpl
${PLIST.WEBUI}tinderbox/scripts/webui/templates/paefchen/display_login.tpl
${PLIST.WEBUI}tinderbox/scripts/webui/templates/paefchen/display_build_groups.tpl
${PLIST.WEBUI}tinderbox/scripts/webui/templates/paefchen/display_markup_log.tpl
${PLIST.WEBUI}tinderbox/scripts/webui/templates/paefchen/failed_buildports.tpl
${PLIST.WEBUI}tinderbox/scripts/webui/templates/paefchen/footer.inc.tpl
${PLIST.WEBUI}tinderbox/scripts/webui/templates/paefchen/header.inc.tpl
${PLIST.WEBUI}tinderbox/scripts/webui/templates/paefchen/images/hdr_fill.png
${PLIST.WEBUI}tinderbox/scripts/webui/templates/paefchen/latest_buildports.tpl
${PLIST.WEBUI}tinderbox/scripts/webui/templates/paefchen/list_buildports.tpl
${PLIST.WEBUI}tinderbox/scripts/webui/templates/paefchen/list_builds.tpl
${PLIST.WEBUI}tinderbox/scripts/webui/templates/paefchen/list_failure_reasons.tpl
${PLIST.WEBUI}tinderbox/scripts/webui/templates/paefchen/list_tinderd_queue.tpl
${PLIST.WEBUI}tinderbox/scripts/webui/templates/paefchen/messages.inc
${PLIST.WEBUI}tinderbox/scripts/webui/templates/paefchen/please_login.tpl
${PLIST.WEBUI}tinderbox/scripts/webui/templates/paefchen/rss.tpl
${PLIST.WEBUI}tinderbox/scripts/webui/templates/paefchen/tinderstyle.css
${PLIST.WEBUI}tinderbox/scripts/webui/templates/paefchen/tinderstyle.js
${PLIST.WEBUI}tinderbox/scripts/webui/templates/paefchen/user_admin.tpl
${PLIST.WEBUI}tinderbox/scripts/webui/templates/paefchen/user_permissions.tpl
${PLIST.WEBUI}tinderbox/scripts/webui/templates/paefchen/user_properties.tpl
