# www/shibboleth-sp/PLIST
@comment $NetBSD: PLIST,v 1.2 2011/07/28 23:16:24 pettai Exp $
bin/mdquery
bin/resolvertest
bin/shib-metagen
include/shibsp/AbstractSPRequest.h
include/shibsp/AccessControl.h
include/shibsp/Application.h
include/shibsp/GSSRequest.h
include/shibsp/RequestMapper.h
include/shibsp/SPConfig.h
include/shibsp/SPRequest.h
include/shibsp/ServiceProvider.h
include/shibsp/SessionCache.h
include/shibsp/SessionCacheEx.h
include/shibsp/TransactionLog.h
include/shibsp/attribute/Attribute.h
include/shibsp/attribute/AttributeDecoder.h
include/shibsp/attribute/ExtensibleAttribute.h
include/shibsp/attribute/NameIDAttribute.h
include/shibsp/attribute/ScopedAttribute.h
include/shibsp/attribute/SimpleAttribute.h
include/shibsp/attribute/XMLAttribute.h
include/shibsp/attribute/filtering/AttributeFilter.h
include/shibsp/attribute/filtering/BasicFilteringContext.h
include/shibsp/attribute/filtering/FilterPolicyContext.h
include/shibsp/attribute/filtering/FilteringContext.h
include/shibsp/attribute/filtering/MatchFunctor.h
include/shibsp/attribute/resolver/AttributeExtractor.h
include/shibsp/attribute/resolver/AttributeResolver.h
include/shibsp/attribute/resolver/ResolutionContext.h
include/shibsp/base.h
include/shibsp/binding/ArtifactResolver.h
include/shibsp/binding/ProtocolProvider.h
include/shibsp/binding/SOAPClient.h
include/shibsp/config_pub.h
include/shibsp/exceptions.h
include/shibsp/handler/AbstractHandler.h
include/shibsp/handler/AssertionConsumerService.h
include/shibsp/handler/Handler.h
include/shibsp/handler/LogoutHandler.h
include/shibsp/handler/LogoutInitiator.h
include/shibsp/handler/RemotedHandler.h
include/shibsp/handler/SessionInitiator.h
include/shibsp/lite/CommonDomainCookie.h
include/shibsp/lite/SAMLConstants.h
include/shibsp/metadata/MetadataExt.h
include/shibsp/metadata/MetadataProviderCriteria.h
include/shibsp/paths.h
include/shibsp/remoting/ListenerService.h
include/shibsp/remoting/ddf.h
include/shibsp/security/PKIXTrustEngine.h
include/shibsp/security/SecurityPolicy.h
include/shibsp/security/SecurityPolicyProvider.h
include/shibsp/util/CGIParser.h
include/shibsp/util/DOMPropertySet.h
include/shibsp/util/PropertySet.h
include/shibsp/util/SPConstants.h
include/shibsp/util/TemplateParameters.h
include/shibsp/version.h
lib/libshibsp-lite.a
lib/libshibsp-lite.so
lib/libshibsp-lite.so.5
lib/libshibsp-lite.so.5.0.3
lib/libshibsp.a
lib/libshibsp.so
lib/libshibsp.so.5
lib/libshibsp.so.5.0.3
lib/shibboleth/adfs-lite.la
lib/shibboleth/adfs.la
lib/shibboleth/mod_shib_22.la
sbin/shib-keygen
sbin/shibd
share/doc/shibboleth-${PKGVERSION}/CREDITS.txt
share/doc/shibboleth-${PKGVERSION}/FASTCGI.LICENSE
share/doc/shibboleth-${PKGVERSION}/LICENSE.txt
share/doc/shibboleth-${PKGVERSION}/LOG4CPP.LICENSE
share/doc/shibboleth-${PKGVERSION}/NOTICE.txt
share/doc/shibboleth-${PKGVERSION}/OPENSSL.LICENSE
share/doc/shibboleth-${PKGVERSION}/README.txt
share/doc/shibboleth-${PKGVERSION}/RELEASE.txt
share/doc/shibboleth-${PKGVERSION}/logo.jpg
share/doc/shibboleth-${PKGVERSION}/main.css
share/examples/rc.d/shibd
share/examples/shibboleth/accessError.html
share/examples/shibboleth/accessError.html.dist
share/examples/shibboleth/apache.config
share/examples/shibboleth/apache2.config
share/examples/shibboleth/apache22.config
share/examples/shibboleth/attribute-map.xml
share/examples/shibboleth/attribute-map.xml.dist
share/examples/shibboleth/attribute-policy.xml
share/examples/shibboleth/attribute-policy.xml.dist
share/examples/shibboleth/bindingTemplate.html
share/examples/shibboleth/bindingTemplate.html.dist
share/examples/shibboleth/console.logger
share/examples/shibboleth/console.logger.dist
share/examples/shibboleth/discoveryTemplate.html
share/examples/shibboleth/discoveryTemplate.html.dist
share/examples/shibboleth/example-metadata.xml
share/examples/shibboleth/example-metadata.xml.dist
share/examples/shibboleth/example-shibboleth2.xml
share/examples/shibboleth/example-shibboleth2.xml.dist
share/examples/shibboleth/globalLogout.html
share/examples/shibboleth/globalLogout.html.dist
share/examples/shibboleth/localLogout.html
share/examples/shibboleth/localLogout.html.dist
share/examples/shibboleth/metadataError.html
share/examples/shibboleth/metadataError.html.dist
share/examples/shibboleth/native.logger
share/examples/shibboleth/native.logger.dist
share/examples/shibboleth/partialLogout.html
share/examples/shibboleth/partialLogout.html.dist
share/examples/shibboleth/postTemplate.html
share/examples/shibboleth/postTemplate.html.dist
share/examples/shibboleth/protocols.xml
share/examples/shibboleth/protocols.xml.dist
share/examples/shibboleth/security-policy.xml
share/examples/shibboleth/security-policy.xml.dist
share/examples/shibboleth/sessionError.html
share/examples/shibboleth/sessionError.html.dist
share/examples/shibboleth/shibboleth2.xml
share/examples/shibboleth/shibboleth2.xml.dist
share/examples/shibboleth/shibd-osx.plist
share/examples/shibboleth/shibd.logger
share/examples/shibboleth/shibd.logger.dist
share/examples/shibboleth/sslError.html
share/examples/shibboleth/sslError.html.dist
share/examples/shibboleth/syslog.logger
share/examples/shibboleth/syslog.logger.dist
share/examples/shibboleth/upgrade.xsl
share/xml/shibboleth/WS-Trust.xsd
share/xml/shibboleth/catalog.xml
share/xml/shibboleth/shibboleth-2.0-afp-mf-basic.xsd
share/xml/shibboleth/shibboleth-2.0-afp-mf-saml.xsd
share/xml/shibboleth/shibboleth-2.0-afp.xsd
share/xml/shibboleth/shibboleth-2.0-attribute-map.xsd
share/xml/shibboleth/shibboleth-2.0-native-sp-config.xsd
share/xml/shibboleth/shibboleth-2.0-native-sp-protocols.xsd
share/xml/shibboleth/shibboleth-2.0-sp-notify.xsd
share/xml/shibboleth/shibboleth-metadata-1.0.xsd
share/xml/shibboleth/shibboleth.xsd
@pkgdir share/doc/shibboleth-${PKGVERSION}/api
@pkgdir etc/shibboleth
