# sysutils/bcmfw/PLIST
@comment $NetBSD: PLIST,v 1.2 2009/06/14 18:16:01 joerg Exp $
man/man8/bcmfw.8
sbin/bcmfw
share/examples/rc.d/bcmfw
${FIRMWARE_DIR}/BCM-LEGAL.txt
${FIRMWARE_DIR}/BCM2033-FW.bin
${FIRMWARE_DIR}/BCM2033-MD.hex
