# security/ssh/PLIST
@comment $NetBSD: PLIST,v 1.3 2002/06/26 10:30:01 seb Exp $
bin/ssh
bin/ssh1
bin/scp
bin/scp1
bin/slogin
bin/ssh-add
bin/ssh-add1
bin/ssh-agent
bin/ssh-agent1
bin/ssh-keygen
bin/ssh-keygen1
bin/make-ssh-known-hosts
bin/make-ssh-known-hosts1
etc/rc.d/sshd
man/man1/make-ssh-known-hosts1.1
man/man1/make-ssh-known-hosts.1
man/man1/scp.1
man/man1/scp1.1
man/man1/ssh-add.1
man/man1/ssh-add1.1
man/man1/ssh-agent.1
man/man1/ssh-agent1.1
man/man1/ssh-keygen.1
man/man1/ssh-keygen1.1
man/man1/ssh.1
man/man1/ssh1.1
man/man1/slogin.1
man/man1/slogin1.1
man/man8/sshd.8
man/man8/sshd1.8
sbin/sshd
sbin/sshd1
share/examples/ssh/ssh_config
share/examples/ssh/sshd_config
@exec if [ ! -d @SSH_CONF_DIR@ ]; then echo "Creating directory @SSH_CONF_DIR@ for ssh config files.." ; ${MKDIR} @SSH_CONF_DIR@; fi
@exec if [ ! -f @SSH_CONF_DIR@/ssh_config ]; then echo "Installing example ssh_config in @SSH_CONF_DIR@.." ; ${INSTALL} -c -o root -g ${ROOT_GROUP} -m 0644 %D/share/examples/ssh/ssh_config @SSH_CONF_DIR@; fi
@exec if [ ! -f @SSH_CONF_DIR@/sshd_config ]; then echo "Installing example sshd_config in @SSH_CONF_DIR@.." ; ${INSTALL} -c -o root -g ${ROOT_GROUP} -m 0644 %D/share/examples/ssh/sshd_config @SSH_CONF_DIR@; fi
@exec if [ ! -f @SSH_CONF_DIR@/ssh_host_key ]; then echo "Generating a secret host key in @SSH_CONF_DIR@.." ; %D/bin/ssh-keygen -N "" -f @SSH_CONF_DIR@/ssh_host_key; fi
@dirrm share/examples/ssh
