# mail/mailgraph/PLIST
@comment $NetBSD: PLIST,v 1.4 2009/06/14 18:04:37 joerg Exp $
libexec/cgi-bin/mailgraph.cgi
sbin/mailgraph
share/doc/mailgraph/CHANGES
share/doc/mailgraph/COPYING
share/doc/mailgraph/README
share/examples/rc.d/mailgraph
