# mail/mailgraph/PLIST
@comment $NetBSD: PLIST,v 1.2 2005/05/02 20:34:00 reed Exp $
sbin/mailgraph
${CGIDIR}/mailgraph.cgi
${DOCDIR}/CHANGES
${DOCDIR}/COPYING
${DOCDIR}/README
share/examples/rc.d/mailgraph
