# comms/conserver8/PLIST
@comment $NetBSD: PLIST,v 1.3 2006/05/11 21:42:43 seb Exp $
bin/console
lib/conserver/convert
man/man1/console.1
man/man5/conserver.cf.5
man/man5/conserver.passwd.5
man/man8/conserver.8
sbin/conserver
share/doc/conserver/INSTALL
share/doc/conserver/PROTOCOL
share/doc/conserver/README
share/examples/conserver/conserver.cf
share/examples/conserver/conserver.passwd
share/examples/conserver/conserver.rc
share/examples/rc.d/conserver8
@dirrm share/examples/conserver
@dirrm share/doc/conserver
@dirrm lib/conserver
