# chat/silc-server/PLIST
@comment $NetBSD: PLIST,v 1.17 2005/05/02 20:33:58 reed Exp $
lib/silc-server/libsilc.la
lib/silc-server/modules/aes.sim.so
lib/silc-server/modules/blowfish.sim.so
lib/silc-server/modules/cast.sim.so
lib/silc-server/modules/md5.sim.so
lib/silc-server/modules/none.sim.so
lib/silc-server/modules/rc5.sim.so
lib/silc-server/modules/rsa.sim.so
lib/silc-server/modules/sha1.sim.so
lib/silc-server/modules/twofish.sim.so
man/man5/silcd.conf.5
man/man8/silcd.8
sbin/silcd
share/doc/silc-server/CHANGES
share/doc/silc-server/COPYING
share/doc/silc-server/CREDITS
share/doc/silc-server/CodingStyle
share/doc/silc-server/FAQ
share/doc/silc-server/INSTALL
share/doc/silc-server/README
share/doc/silc-server/TODO
share/examples/rc.d/silcd
share/examples/silc-server/README
share/examples/silc-server/cell_backup.conf
share/examples/silc-server/cell_router.conf
share/examples/silc-server/cell_server1.conf
share/examples/silc-server/cell_server2.conf
share/examples/silc-server/motd.txt.default
share/examples/silc-server/silcalgs.conf.default
share/examples/silc-server/silcd.conf.default
share/examples/silc-server/silcd.prv
share/examples/silc-server/silcd.pub
@dirrm share/examples/silc-server
@dirrm share/doc/silc-server
@dirrm lib/silc-server/modules
@dirrm lib/silc-server
