# mail/jchkmail/PLIST
@comment $NetBSD: PLIST,v 1.3 2009/06/14 18:04:35 joerg Exp $
bin/j-scanfile
bin/j-printstats
sbin/j-chkmail
man/man1/j-printstats.1
man/man1/j-scanfile.1
man/man4/j-chkmail.cf.4
man/man8/j-chkmail.8
share/examples/jchkmail/j-chkmail.cf
share/examples/jchkmail/j-local-users
share/examples/jchkmail/j-error-msg
share/examples/jchkmail/j-nets
share/examples/jchkmail/j-host-access
share/examples/jchkmail/j-user-access
share/examples/jchkmail/j-regex
share/examples/rc.d/jchkmail
