# comms/fidogate/pkg/PLIST
@comment $NetBSD: PLIST,v 1.4 1999/02/11 19:17:28 tv Exp $
bin/pktdebug
bin/pktmore
bin/sumcrc
etc/rc.d/fidogate.sh
@unexec install-info --delete %D/info/fidogate.info %D/info/dir
info/fidogate.info
info/fidogate.info-1
info/fidogate.info-2
@exec install-info %D/info/fidogate.info %D/info/dir
lib/fidogate/ffx
lib/fidogate/ffxbatch
lib/fidogate/ffxmail
lib/fidogate/ffxqt
lib/fidogate/ffxrmail
lib/fidogate/ffxrun
lib/fidogate/ftn2ftn
lib/fidogate/ftn2rfc
lib/fidogate/ftnaf
lib/fidogate/ftnafmail
lib/fidogate/ftnbsy
lib/fidogate/ftnconfig
lib/fidogate/ftnexpire
lib/fidogate/ftnfattach
lib/fidogate/ftnflo
lib/fidogate/ftnhatch
lib/fidogate/ftnin
lib/fidogate/ftninpost
lib/fidogate/ftninrecomb
lib/fidogate/ftnlock
lib/fidogate/ftnlog
lib/fidogate/ftnpack
lib/fidogate/ftnroute
lib/fidogate/ftnseq
lib/fidogate/ftntick
lib/fidogate/ftntoss
lib/fidogate/magic
lib/fidogate/rfc2ftn
lib/fidogate/rundf
lib/fidogate/runin
lib/fidogate/runout
lib/fidogate/runtoss
lib/fidogate/rununpack
lib/libfidogate.so.0.0
share/fidogate/inn/send-ffx
share/fidogate/inn/send-fidogate
share/fidogate/inn/viaffx
share/fidogate/nodelist/latest
share/fidogate/nodelist/nl-2hosts
share/fidogate/nodelist/nl-2routing
share/fidogate/nodelist/nl-autoupd
share/fidogate/nodelist/nl-check
share/fidogate/nodelist/nl-del
share/fidogate/nodelist/nl-diff
share/fidogate/nodelist/nl-html
share/fidogate/nodelist/nl-isdn
share/fidogate/nodelist/nl-print
share/fidogate/nodelist/template.html
share/fidogate/scripts/out-freq
share/fidogate/scripts/out-ls
share/fidogate/scripts/out-manip
share/fidogate/scripts/out-rm0
share/fidogate/scripts/out-rmbsy
share/fidogate/scripts/run1
share/fidogate/scripts/run2
share/fidogate/scripts/runffx
share/fidogate/scripts/runpoll
share/fidogate/sendmail/ffx.m4
share/fidogate/sendmail/ftn.m4
@dirrm lib/fidogate
@dirrm share/fidogate/inn
@dirrm share/fidogate/nodelist
@dirrm share/fidogate/scripts
@dirrm share/fidogate/sendmail
@dirrm share/fidogate
@exec install -d -m 755 -o news -g news /var/log/fido /var/spool/fido
@exec for dir in in in/tmpmail in/tmpnews in/bad insecure out outpkt outpkt/mail locks seq toss toss/in toss/bad toss/tmp toss/out toss/pack bt bt/in bt/uuin bt/pin; do install -d -m 755 -o news -g news /var/spool/fido/$dir; done
