# comms/hylafax/PLIST
@comment $NetBSD: PLIST,v 1.12 2009/06/14 17:38:40 joerg Exp $
bin/faxalter
bin/faxcover
bin/faxmail
bin/faxrm
bin/faxstat
bin/sendfax
bin/sendpage
libdata/hylafax/faxcover.ps
libdata/hylafax/faxcover_example_sgi.ps
libdata/hylafax/faxmail.ps
libdata/hylafax/hfaxd.conf
libdata/hylafax/pagesizes
libdata/hylafax/typerules
libexec/faxgetty
libexec/faxsend
libexec/hfaxd
libexec/hylafax/archive
libexec/hylafax/faxrcvd
libexec/hylafax/mkcover
libexec/hylafax/notify
libexec/hylafax/notify-4.1
libexec/hylafax/notify.awk
libexec/hylafax/pcl2fax
libexec/hylafax/pdf2fax.gs
libexec/hylafax/pollrcvd
libexec/hylafax/ps2fax.dps
libexec/hylafax/ps2fax.gs
libexec/hylafax/ps2fax.imp
libexec/hylafax/tiff2fax
libexec/hylafax/tiff2pdf
libexec/hylafax/wedged
libexec/lockname
libexec/ondelay
libexec/pagesend
libexec/textfmt
man/man1/edit-faxcover.1
man/man1/faxalter.1
man/man1/faxcover.1
man/man1/faxmail.1
man/man1/faxrm.1
man/man1/faxstat.1
man/man1/hylafax-client.1
man/man1/sendfax.1
man/man1/sendpage.1
man/man1/sgi2fax.1
man/man1/textfmt.1
man/man5/callid.5
man/man5/destctrls.5
man/man5/dialrules.5
man/man5/doneq.5
man/man5/hosts.hfaxd.5
man/man5/hylafax-config.5
man/man5/hylafax-info.5
man/man5/hylafax-log.5
man/man5/hylafax-server.5
man/man5/hylafax-shutdown.5
man/man5/pagermap.5
man/man5/pagesizes.5
man/man5/recvq.5
man/man5/sendq.5
man/man5/status.5
man/man5/tsi.5
man/man5/typerules.5
man/man5/xferfaxlog.5
man/man8/choptest.8
man/man8/cqtest.8
man/man8/dialtest.8
man/man8/faxabort.8
man/man8/faxaddmodem.8
man/man8/faxadduser.8
man/man8/faxanswer.8
man/man8/faxconfig.8
man/man8/faxcron.8
man/man8/faxdeluser.8
man/man8/faxgetty.8
man/man8/faxinfo.8
man/man8/faxlock.8
man/man8/faxmodem.8
man/man8/faxq.8
man/man8/faxqclean.8
man/man8/faxquit.8
man/man8/faxrcvd.8
man/man8/faxsend.8
man/man8/faxsetup.8
man/man8/faxstate.8
man/man8/faxwatch.8
man/man8/hfaxd.8
man/man8/mkcover.8
man/man8/notify.8
man/man8/pagesend.8
man/man8/pdf2fax.8
man/man8/pollrcvd.8
man/man8/ps2fax.8
man/man8/recvstats.8
man/man8/tagtest.8
man/man8/tiff2fax.8
man/man8/tiffcheck.8
man/man8/tsitest.8
man/man8/wedged.8
man/man8/xferfaxstats.8
sbin/choptest
sbin/cqtest
sbin/dialtest
sbin/edit-faxcover
sbin/faxabort
sbin/faxaddmodem
sbin/faxadduser
sbin/faxanswer
sbin/faxconfig
sbin/faxcron
sbin/faxdeluser
sbin/faxinfo
sbin/faxlock
sbin/faxmodem
sbin/faxmsg
sbin/faxq
sbin/faxqclean
sbin/faxquit
sbin/faxsetup
sbin/faxsetup.bsdi
sbin/faxsetup.irix
sbin/faxsetup.linux
sbin/faxstate
sbin/faxwatch
sbin/hylafax
sbin/probemodem
sbin/recvstats
sbin/tagtest
sbin/tiffcheck
sbin/tsitest
sbin/typetest
sbin/xferfaxstats
share/doc/hylafax/COPYRIGHT
share/examples/rc.d/hylafax
share/hylafax/config_template/agere-mp-2
share/hylafax/config_template/agere-mp-21
share/hylafax/config_template/att-dataport
share/hylafax/config_template/att-dataport-2
share/hylafax/config_template/att-dataport-2.0
share/hylafax/config_template/att-dpexpress-2
share/hylafax/config_template/cirrus
share/hylafax/config_template/class1
share/hylafax/config_template/class1.0
share/hylafax/config_template/class2
share/hylafax/config_template/class2.0
share/hylafax/config_template/class2.1
share/hylafax/config_template/conexant-1.0
share/hylafax/config_template/davicom-2
share/hylafax/config_template/digi
share/hylafax/config_template/digi-2
share/hylafax/config_template/digi-20
share/hylafax/config_template/dsi-scout+
share/hylafax/config_template/e-tech
share/hylafax/config_template/eicon
share/hylafax/config_template/eicon-2
share/hylafax/config_template/elsa-microlink-56k-basic-2.0
share/hylafax/config_template/everex-2496d
share/hylafax/config_template/exar
share/hylafax/config_template/gvc-288
share/hylafax/config_template/gvc-288-2
share/hylafax/config_template/hayes
share/hylafax/config_template/intel-400e
share/hylafax/config_template/lucent
share/hylafax/config_template/lucent-isa
share/hylafax/config_template/lucent-mt-10
share/hylafax/config_template/lucent-mt-2
share/hylafax/config_template/lucent-mt-20
share/hylafax/config_template/lucent-mt-21
share/hylafax/config_template/moto-288
share/hylafax/config_template/mt-1432
share/hylafax/config_template/nuvo-voyager
share/hylafax/config_template/ppi-pm14400fxmt
share/hylafax/config_template/ppi-pm14400fxsa
share/hylafax/config_template/prometheus
share/hylafax/config_template/rc144ac
share/hylafax/config_template/rc144ac-1
share/hylafax/config_template/rc288dpi
share/hylafax/config_template/rc288dpi-1
share/hylafax/config_template/rc32acl
share/hylafax/config_template/rc32acl-1
share/hylafax/config_template/rockwell-k56
share/hylafax/config_template/rockwell-rc
share/hylafax/config_template/sierra
share/hylafax/config_template/skel
share/hylafax/config_template/t38modem
share/hylafax/config_template/telebit-qblazer
share/hylafax/config_template/telebit-t3000
share/hylafax/config_template/telebit-wb
share/hylafax/config_template/topic
share/hylafax/config_template/tricom2842
share/hylafax/config_template/uds-fastalk
share/hylafax/config_template/umc92144
share/hylafax/config_template/umc9624
share/hylafax/config_template/usr-2.0
share/hylafax/config_template/usr-rts
share/hylafax/config_template/usr-xon
share/hylafax/config_template/zyxel-1496e
share/hylafax/config_template/zyxel-1496e-1
share/hylafax/config_template/zyxel-1496e-2.0
share/hylafax/config_template/zyxel-2864
share/hylafax/etc_template/cover.templ
share/hylafax/etc_template/dialrules
share/hylafax/etc_template/dialrules.europe
share/hylafax/etc_template/dialrules.sf-ba
share/hylafax/etc_template/dpsprinter.ps
share/hylafax/etc_template/hosts.hfaxd
share/hylafax/etc_template/lutRS18.pcf
share/hylafax/etc_template/xferfaxlog
