# comms/hylafax/PLIST
@comment $NetBSD: PLIST,v 1.9 2005/05/02 20:33:58 reed Exp $
bin/faxalter
bin/faxcover
bin/faxmail
bin/faxrm
bin/faxstat
bin/sendfax
bin/sendpage
libdata/hylafax/faxcover.ps
libdata/hylafax/faxcover_example_sgi.ps
libdata/hylafax/faxmail.ps
libdata/hylafax/hfaxd.conf
libdata/hylafax/pagesizes
libdata/hylafax/typerules
libexec/faxgetty
libexec/faxsend
libexec/hfaxd
libexec/hylafax/archive
libexec/hylafax/faxrcvd
libexec/hylafax/mkcover
libexec/hylafax/notify
libexec/hylafax/notify-4.1
libexec/hylafax/pcl2fax
libexec/hylafax/pdf2fax.gs
libexec/hylafax/pollrcvd
libexec/hylafax/ps2fax.dps
libexec/hylafax/ps2fax.gs
libexec/hylafax/ps2fax.imp
libexec/hylafax/tiff2fax
libexec/hylafax/tiff2pdf
libexec/hylafax/wedged
libexec/lockname
libexec/ondelay
libexec/pagesend
libexec/textfmt
man/cat1/edit-faxcover.0
man/cat1/faxalter.0
man/cat1/faxcover.0
man/cat1/faxmail.0
man/cat1/faxrm.0
man/cat1/faxstat.0
man/cat1/hylafax-client.0
man/cat1/sendfax.0
man/cat1/sendpage.0
man/cat1/sgi2fax.0
man/cat1/textfmt.0
man/cat5/cid.0
man/cat5/destctrls.0
man/cat5/dialrules.0
man/cat5/doneq.0
man/cat5/hosts.hfaxd.0
man/cat5/hylafax-config.0
man/cat5/hylafax-info.0
man/cat5/hylafax-log.0
man/cat5/hylafax-server.0
man/cat5/hylafax-shutdown.0
man/cat5/pagermap.0
man/cat5/pagesizes.0
man/cat5/recvq.0
man/cat5/sendq.0
man/cat5/status.0
man/cat5/tsi.0
man/cat5/typerules.0
man/cat5/xferfaxlog.0
man/cat8/choptest.0
man/cat8/cqtest.0
man/cat8/dialtest.0
man/cat8/faxabort.0
man/cat8/faxaddmodem.0
man/cat8/faxadduser.0
man/cat8/faxanswer.0
man/cat8/faxconfig.0
man/cat8/faxcron.0
man/cat8/faxdeluser.0
man/cat8/faxgetty.0
man/cat8/faxinfo.0
man/cat8/faxlock.0
man/cat8/faxmodem.0
man/cat8/faxq.0
man/cat8/faxqclean.0
man/cat8/faxquit.0
man/cat8/faxrcvd.0
man/cat8/faxsend.0
man/cat8/faxsetup.0
man/cat8/faxstate.0
man/cat8/faxwatch.0
man/cat8/hfaxd.0
man/cat8/mkcover.0
man/cat8/notify.0
man/cat8/pagesend.0
man/cat8/pdf2fax.0
man/cat8/pollrcvd.0
man/cat8/ps2fax.0
man/cat8/recvstats.0
man/cat8/tagtest.0
man/cat8/tiff2fax.0
man/cat8/tiffcheck.0
man/cat8/tsitest.0
man/cat8/wedged.0
man/cat8/xferfaxstats.0
sbin/choptest
sbin/cqtest
sbin/dialtest
sbin/edit-faxcover
sbin/faxabort
sbin/faxaddmodem
sbin/faxadduser
sbin/faxanswer
sbin/faxconfig
sbin/faxcron
sbin/faxdeluser
sbin/faxinfo
sbin/faxlock
sbin/faxmodem
sbin/faxmsg
sbin/faxq
sbin/faxqclean
sbin/faxquit
sbin/faxsetup
sbin/faxsetup.bsdi
sbin/faxsetup.irix
sbin/faxsetup.linux
sbin/faxstate
sbin/faxwatch
sbin/hylafax
sbin/probemodem
sbin/recvstats
sbin/tagtest
sbin/tiffcheck
sbin/tsitest
sbin/typetest
sbin/xferfaxstats
share/doc/hylafax/COPYRIGHT
share/examples/rc.d/hylafax
share/hylafax/config_template/agere-mp-2
share/hylafax/config_template/agere-mp-21
share/hylafax/config_template/att-dataport
share/hylafax/config_template/att-dataport-2
share/hylafax/config_template/att-dataport-2.0
share/hylafax/config_template/att-dpexpress-2
share/hylafax/config_template/cirrus
share/hylafax/config_template/class1
share/hylafax/config_template/class1.0
share/hylafax/config_template/class2
share/hylafax/config_template/class2.0
share/hylafax/config_template/class2.1
share/hylafax/config_template/davicom-2
share/hylafax/config_template/digi
share/hylafax/config_template/digi-2
share/hylafax/config_template/digi-20
share/hylafax/config_template/dsi-scout+
share/hylafax/config_template/e-tech
share/hylafax/config_template/eicon
share/hylafax/config_template/eicon-2
share/hylafax/config_template/elsa-microlink-56k-basic-2.0
share/hylafax/config_template/everex-2496d
share/hylafax/config_template/exar
share/hylafax/config_template/gvc-288
share/hylafax/config_template/gvc-288-2
share/hylafax/config_template/hayes
share/hylafax/config_template/intel-400e
share/hylafax/config_template/lucent
share/hylafax/config_template/lucent-isa
share/hylafax/config_template/lucent-mt-10
share/hylafax/config_template/lucent-mt-2
share/hylafax/config_template/lucent-mt-20
share/hylafax/config_template/lucent-mt-21
share/hylafax/config_template/moto-288
share/hylafax/config_template/mt-1432
share/hylafax/config_template/nuvo-voyager
share/hylafax/config_template/ppi-pm14400fxmt
share/hylafax/config_template/ppi-pm14400fxsa
share/hylafax/config_template/prometheus
share/hylafax/config_template/rc144ac
share/hylafax/config_template/rc144ac-1
share/hylafax/config_template/rc288dpi
share/hylafax/config_template/rc288dpi-1
share/hylafax/config_template/rc32acl
share/hylafax/config_template/rc32acl-1
share/hylafax/config_template/rockwell-k56
share/hylafax/config_template/rockwell-rc
share/hylafax/config_template/sierra
share/hylafax/config_template/skel
share/hylafax/config_template/t38modem
share/hylafax/config_template/telebit-qblazer
share/hylafax/config_template/telebit-t3000
share/hylafax/config_template/telebit-wb
share/hylafax/config_template/topic
share/hylafax/config_template/tricom2842
share/hylafax/config_template/uds-fastalk
share/hylafax/config_template/umc92144
share/hylafax/config_template/umc9624
share/hylafax/config_template/usr-2.0
share/hylafax/config_template/usr-rts
share/hylafax/config_template/usr-xon
share/hylafax/config_template/zyxel-1496e
share/hylafax/config_template/zyxel-1496e-1
share/hylafax/config_template/zyxel-1496e-2.0
share/hylafax/config_template/zyxel-2864
share/hylafax/etc_template/cover.templ
share/hylafax/etc_template/dialrules
share/hylafax/etc_template/dialrules.europe
share/hylafax/etc_template/dialrules.sf-ba
share/hylafax/etc_template/dpsprinter.ps
share/hylafax/etc_template/hosts.hfaxd
share/hylafax/etc_template/lutRS18.pcf
share/hylafax/etc_template/xferfaxlog
@dirrm share/hylafax/etc_template
@dirrm share/hylafax/config_template
@dirrm share/hylafax
@dirrm share/doc/hylafax
@dirrm libexec/hylafax
@dirrm libdata/hylafax
