# security/zoneminder/PLIST
@comment $NetBSD: PLIST,v 1.3 2014/05/06 15:07:40 jperkin Exp $
bin/zma
bin/zmaudit.pl
bin/zmc
bin/zmcontrol.pl
bin/zmdc.pl
bin/zmf
bin/zmfilter.pl
bin/zmfix
bin/zmpkg.pl
bin/zmstreamer
bin/zmtrack.pl
bin/zmtrigger.pl
bin/zmu
bin/zmupdate.pl
bin/zmvideo.pl
bin/zmwatch.pl
bin/zmx10.pl
${ZM_HTTPD_CGIBIN}/nph-zms
${ZM_HTTPD_CGIBIN}/zms
${PLIST.rcd}share/examples/rc.d/zoneminder
share/examples/zoneminder/apache/zoneminder.conf
share/examples/zoneminder/config/zm.conf
share/zoneminder/db/zm_create.sql
share/zoneminder/db/zm_update-0.0.1.sql
share/zoneminder/db/zm_update-0.9.10.sql
share/zoneminder/db/zm_update-0.9.11.sql
share/zoneminder/db/zm_update-0.9.12.sql
share/zoneminder/db/zm_update-0.9.13.sql
share/zoneminder/db/zm_update-0.9.15.sql
share/zoneminder/db/zm_update-0.9.16.sql
share/zoneminder/db/zm_update-0.9.7.sql
share/zoneminder/db/zm_update-0.9.8.sql
share/zoneminder/db/zm_update-0.9.9.sql
share/zoneminder/db/zm_update-1.17.1.sql
share/zoneminder/db/zm_update-1.17.2.sql
share/zoneminder/db/zm_update-1.18.0.sql
share/zoneminder/db/zm_update-1.18.1.sql
share/zoneminder/db/zm_update-1.19.0.sql
share/zoneminder/db/zm_update-1.19.1.sql
share/zoneminder/db/zm_update-1.19.2.sql
share/zoneminder/db/zm_update-1.19.3.sql
share/zoneminder/db/zm_update-1.19.4.sql
share/zoneminder/db/zm_update-1.19.5.sql
share/zoneminder/db/zm_update-1.20.0.sql
share/zoneminder/db/zm_update-1.20.1.sql
share/zoneminder/db/zm_update-1.21.0.sql
share/zoneminder/db/zm_update-1.21.1.sql
share/zoneminder/db/zm_update-1.21.2.sql
share/zoneminder/db/zm_update-1.21.3.sql
share/zoneminder/db/zm_update-1.21.4.sql
share/zoneminder/db/zm_update-1.22.0.sql
share/zoneminder/db/zm_update-1.22.1.sql
share/zoneminder/db/zm_update-1.22.2.sql
share/zoneminder/db/zm_update-1.22.3.sql
share/zoneminder/db/zm_update-1.23.0.sql
share/zoneminder/db/zm_update-1.23.1.sql
share/zoneminder/db/zm_update-1.23.2.sql
share/zoneminder/db/zm_update-1.23.3.sql
share/zoneminder/db/zm_update-1.24.0.sql
share/zoneminder/db/zm_update-1.24.1.sql
share/zoneminder/db/zm_update-1.24.2.sql
share/zoneminder/db/zm_update-1.24.3.sql
share/zoneminder/db/zm_update-1.24.4.sql
${ZM_HTTPD_ROOT}/ajax/alarm.php
${ZM_HTTPD_ROOT}/ajax/control.php
${ZM_HTTPD_ROOT}/ajax/event.php
${ZM_HTTPD_ROOT}/ajax/log.php
${ZM_HTTPD_ROOT}/ajax/status.php
${ZM_HTTPD_ROOT}/ajax/stream.php
${ZM_HTTPD_ROOT}/ajax/zone.php
${ZM_HTTPD_ROOT}/css/overlay.css
${ZM_HTTPD_ROOT}/css/reset.css
${ZM_HTTPD_ROOT}/css/spinner.css
${ZM_HTTPD_ROOT}/events
${ZM_HTTPD_ROOT}/graphics/favicon.ico
${ZM_HTTPD_ROOT}/graphics/spinner.gif
${ZM_HTTPD_ROOT}/graphics/transparent.gif
${ZM_HTTPD_ROOT}/images
${ZM_HTTPD_ROOT}/includes/actions.php
${ZM_HTTPD_ROOT}/includes/config.php
${ZM_HTTPD_ROOT}/includes/control_functions.php
${ZM_HTTPD_ROOT}/includes/database.php
${ZM_HTTPD_ROOT}/includes/functions.php
${ZM_HTTPD_ROOT}/includes/lang.php
${ZM_HTTPD_ROOT}/includes/logger.php
${ZM_HTTPD_ROOT}/index.php
${ZM_HTTPD_ROOT}/js/logger.js
${ZM_HTTPD_ROOT}/js/mootools.ext.js
${ZM_HTTPD_ROOT}/js/overlay.js
${ZM_HTTPD_ROOT}/lang/big5_big5.php
${ZM_HTTPD_ROOT}/lang/cn_zh.php
${ZM_HTTPD_ROOT}/lang/cs_cz.php
${ZM_HTTPD_ROOT}/lang/de_de.php
${ZM_HTTPD_ROOT}/lang/dk_dk.php
${ZM_HTTPD_ROOT}/lang/en_gb.php
${ZM_HTTPD_ROOT}/lang/en_us.php
${ZM_HTTPD_ROOT}/lang/es_ar.php
${ZM_HTTPD_ROOT}/lang/et_ee.php
${ZM_HTTPD_ROOT}/lang/fr_fr.php
${ZM_HTTPD_ROOT}/lang/he_il.php
${ZM_HTTPD_ROOT}/lang/hu_hu.php
${ZM_HTTPD_ROOT}/lang/it_it.php
${ZM_HTTPD_ROOT}/lang/ja_jp.php
${ZM_HTTPD_ROOT}/lang/nl_nl.php
${ZM_HTTPD_ROOT}/lang/pl_pl.php
${ZM_HTTPD_ROOT}/lang/pt_br.php
${ZM_HTTPD_ROOT}/lang/ro_ro.php
${ZM_HTTPD_ROOT}/lang/ru_ru.php
${ZM_HTTPD_ROOT}/lang/se_se.php
${ZM_HTTPD_ROOT}/skins/classic/css/control.css
${ZM_HTTPD_ROOT}/skins/classic/css/export.css
${ZM_HTTPD_ROOT}/skins/classic/css/skin.css
${ZM_HTTPD_ROOT}/skins/classic/graphics/arrow-d.gif
${ZM_HTTPD_ROOT}/skins/classic/graphics/arrow-dl.gif
${ZM_HTTPD_ROOT}/skins/classic/graphics/arrow-dr.gif
${ZM_HTTPD_ROOT}/skins/classic/graphics/arrow-l-d.gif
${ZM_HTTPD_ROOT}/skins/classic/graphics/arrow-l-u.gif
${ZM_HTTPD_ROOT}/skins/classic/graphics/arrow-l.gif
${ZM_HTTPD_ROOT}/skins/classic/graphics/arrow-r.gif
${ZM_HTTPD_ROOT}/skins/classic/graphics/arrow-s-d.gif
${ZM_HTTPD_ROOT}/skins/classic/graphics/arrow-s-u.gif
${ZM_HTTPD_ROOT}/skins/classic/graphics/arrow-u.gif
${ZM_HTTPD_ROOT}/skins/classic/graphics/arrow-ul.gif
${ZM_HTTPD_ROOT}/skins/classic/graphics/arrow-ur.gif
${ZM_HTTPD_ROOT}/skins/classic/graphics/center.gif
${ZM_HTTPD_ROOT}/skins/classic/graphics/point-g.gif
${ZM_HTTPD_ROOT}/skins/classic/graphics/point-o.gif
${ZM_HTTPD_ROOT}/skins/classic/graphics/point-r.gif
${ZM_HTTPD_ROOT}/skins/classic/graphics/seq-d.gif
${ZM_HTTPD_ROOT}/skins/classic/graphics/seq-u.gif
${ZM_HTTPD_ROOT}/skins/classic/includes/config.php
${ZM_HTTPD_ROOT}/skins/classic/includes/control_functions.php
${ZM_HTTPD_ROOT}/skins/classic/includes/export_functions.php
${ZM_HTTPD_ROOT}/skins/classic/includes/functions.php
${ZM_HTTPD_ROOT}/skins/classic/includes/init.php
${ZM_HTTPD_ROOT}/skins/classic/includes/timeline_functions.php
${ZM_HTTPD_ROOT}/skins/classic/js/skin.js
${ZM_HTTPD_ROOT}/skins/classic/js/skin.js.php
${ZM_HTTPD_ROOT}/skins/classic/skin.php
${ZM_HTTPD_ROOT}/skins/classic/views/Makefile.am
${ZM_HTTPD_ROOT}/skins/classic/views/bandwidth.php
${ZM_HTTPD_ROOT}/skins/classic/views/blank.php
${ZM_HTTPD_ROOT}/skins/classic/views/console.php
${ZM_HTTPD_ROOT}/skins/classic/views/control.php
${ZM_HTTPD_ROOT}/skins/classic/views/controlcap.php
${ZM_HTTPD_ROOT}/skins/classic/views/controlcaps.php
${ZM_HTTPD_ROOT}/skins/classic/views/controlpreset.php
${ZM_HTTPD_ROOT}/skins/classic/views/css/console.css
${ZM_HTTPD_ROOT}/skins/classic/views/css/control.css
${ZM_HTTPD_ROOT}/skins/classic/views/css/controlcaps.css
${ZM_HTTPD_ROOT}/skins/classic/views/css/devices.css
${ZM_HTTPD_ROOT}/skins/classic/views/css/event.css
${ZM_HTTPD_ROOT}/skins/classic/views/css/events.css
${ZM_HTTPD_ROOT}/skins/classic/views/css/export.css
${ZM_HTTPD_ROOT}/skins/classic/views/css/filter.css
${ZM_HTTPD_ROOT}/skins/classic/views/css/frame.css
${ZM_HTTPD_ROOT}/skins/classic/views/css/frames.css
${ZM_HTTPD_ROOT}/skins/classic/views/css/groups.css
${ZM_HTTPD_ROOT}/skins/classic/views/css/log.css
${ZM_HTTPD_ROOT}/skins/classic/views/css/monitor.css
${ZM_HTTPD_ROOT}/skins/classic/views/css/montage.css
${ZM_HTTPD_ROOT}/skins/classic/views/css/montage_2wide.css
${ZM_HTTPD_ROOT}/skins/classic/views/css/montage_3wide.css
${ZM_HTTPD_ROOT}/skins/classic/views/css/montage_3wide50enlarge.css
${ZM_HTTPD_ROOT}/skins/classic/views/css/montage_4wide.css
${ZM_HTTPD_ROOT}/skins/classic/views/css/montage_freeform.css
${ZM_HTTPD_ROOT}/skins/classic/views/css/options.css
${ZM_HTTPD_ROOT}/skins/classic/views/css/stats.css
${ZM_HTTPD_ROOT}/skins/classic/views/css/timeline.css
${ZM_HTTPD_ROOT}/skins/classic/views/css/timeline.css.php
${ZM_HTTPD_ROOT}/skins/classic/views/css/video.css
${ZM_HTTPD_ROOT}/skins/classic/views/css/watch.css
${ZM_HTTPD_ROOT}/skins/classic/views/css/zone.css
${ZM_HTTPD_ROOT}/skins/classic/views/cycle.php
${ZM_HTTPD_ROOT}/skins/classic/views/device.php
${ZM_HTTPD_ROOT}/skins/classic/views/devices.php
${ZM_HTTPD_ROOT}/skins/classic/views/donate.php
${ZM_HTTPD_ROOT}/skins/classic/views/error.php
${ZM_HTTPD_ROOT}/skins/classic/views/event.php
${ZM_HTTPD_ROOT}/skins/classic/views/eventdetail.php
${ZM_HTTPD_ROOT}/skins/classic/views/events.php
${ZM_HTTPD_ROOT}/skins/classic/views/export.php
${ZM_HTTPD_ROOT}/skins/classic/views/filter.php
${ZM_HTTPD_ROOT}/skins/classic/views/filtersave.php
${ZM_HTTPD_ROOT}/skins/classic/views/frame.php
${ZM_HTTPD_ROOT}/skins/classic/views/frames.php
${ZM_HTTPD_ROOT}/skins/classic/views/function.php
${ZM_HTTPD_ROOT}/skins/classic/views/group.php
${ZM_HTTPD_ROOT}/skins/classic/views/groups.php
${ZM_HTTPD_ROOT}/skins/classic/views/js/Makefile.am
${ZM_HTTPD_ROOT}/skins/classic/views/js/console.js
${ZM_HTTPD_ROOT}/skins/classic/views/js/console.js.php
${ZM_HTTPD_ROOT}/skins/classic/views/js/control.js
${ZM_HTTPD_ROOT}/skins/classic/views/js/controlpreset.js
${ZM_HTTPD_ROOT}/skins/classic/views/js/controlpreset.js.php
${ZM_HTTPD_ROOT}/skins/classic/views/js/cycle.js
${ZM_HTTPD_ROOT}/skins/classic/views/js/cycle.js.php
${ZM_HTTPD_ROOT}/skins/classic/views/js/devices.js
${ZM_HTTPD_ROOT}/skins/classic/views/js/donate.js
${ZM_HTTPD_ROOT}/skins/classic/views/js/donate.js.php
${ZM_HTTPD_ROOT}/skins/classic/views/js/event.js
${ZM_HTTPD_ROOT}/skins/classic/views/js/event.js.php
${ZM_HTTPD_ROOT}/skins/classic/views/js/events.js
${ZM_HTTPD_ROOT}/skins/classic/views/js/events.js.php
${ZM_HTTPD_ROOT}/skins/classic/views/js/export.js
${ZM_HTTPD_ROOT}/skins/classic/views/js/export.js.php
${ZM_HTTPD_ROOT}/skins/classic/views/js/filter.js
${ZM_HTTPD_ROOT}/skins/classic/views/js/filter.js.php
${ZM_HTTPD_ROOT}/skins/classic/views/js/group.js
${ZM_HTTPD_ROOT}/skins/classic/views/js/groups.js
${ZM_HTTPD_ROOT}/skins/classic/views/js/log.js
${ZM_HTTPD_ROOT}/skins/classic/views/js/login.js
${ZM_HTTPD_ROOT}/skins/classic/views/js/monitor.js
${ZM_HTTPD_ROOT}/skins/classic/views/js/monitor.js.php
${ZM_HTTPD_ROOT}/skins/classic/views/js/monitorpreset.js
${ZM_HTTPD_ROOT}/skins/classic/views/js/monitorprobe.js
${ZM_HTTPD_ROOT}/skins/classic/views/js/montage.js
${ZM_HTTPD_ROOT}/skins/classic/views/js/montage.js.php
${ZM_HTTPD_ROOT}/skins/classic/views/js/options.js.php
${ZM_HTTPD_ROOT}/skins/classic/views/js/postlogin.js
${ZM_HTTPD_ROOT}/skins/classic/views/js/state.js
${ZM_HTTPD_ROOT}/skins/classic/views/js/state.js.php
${ZM_HTTPD_ROOT}/skins/classic/views/js/timeline.js
${ZM_HTTPD_ROOT}/skins/classic/views/js/timeline.js.php
${ZM_HTTPD_ROOT}/skins/classic/views/js/user.js
${ZM_HTTPD_ROOT}/skins/classic/views/js/version.js
${ZM_HTTPD_ROOT}/skins/classic/views/js/version.js.php
${ZM_HTTPD_ROOT}/skins/classic/views/js/video.js
${ZM_HTTPD_ROOT}/skins/classic/views/js/video.js.php
${ZM_HTTPD_ROOT}/skins/classic/views/js/watch.js
${ZM_HTTPD_ROOT}/skins/classic/views/js/watch.js.php
${ZM_HTTPD_ROOT}/skins/classic/views/js/zone.js
${ZM_HTTPD_ROOT}/skins/classic/views/js/zone.js.php
${ZM_HTTPD_ROOT}/skins/classic/views/log.php
${ZM_HTTPD_ROOT}/skins/classic/views/login.php
${ZM_HTTPD_ROOT}/skins/classic/views/logout.php
${ZM_HTTPD_ROOT}/skins/classic/views/monitor.php
${ZM_HTTPD_ROOT}/skins/classic/views/monitorpreset.php
${ZM_HTTPD_ROOT}/skins/classic/views/monitorprobe.php
${ZM_HTTPD_ROOT}/skins/classic/views/montage.php
${ZM_HTTPD_ROOT}/skins/classic/views/none.php
${ZM_HTTPD_ROOT}/skins/classic/views/optionhelp.php
${ZM_HTTPD_ROOT}/skins/classic/views/options.php
${ZM_HTTPD_ROOT}/skins/classic/views/postlogin.php
${ZM_HTTPD_ROOT}/skins/classic/views/settings.php
${ZM_HTTPD_ROOT}/skins/classic/views/state.php
${ZM_HTTPD_ROOT}/skins/classic/views/stats.php
${ZM_HTTPD_ROOT}/skins/classic/views/status.php
${ZM_HTTPD_ROOT}/skins/classic/views/timeline.php
${ZM_HTTPD_ROOT}/skins/classic/views/user.php
${ZM_HTTPD_ROOT}/skins/classic/views/version.php
${ZM_HTTPD_ROOT}/skins/classic/views/video.php
${ZM_HTTPD_ROOT}/skins/classic/views/watch.php
${ZM_HTTPD_ROOT}/skins/classic/views/zone.php
${ZM_HTTPD_ROOT}/skins/classic/views/zones.php
${ZM_HTTPD_ROOT}/skins/mobile/css/skin.css
${ZM_HTTPD_ROOT}/skins/mobile/includes/config.php
${ZM_HTTPD_ROOT}/skins/mobile/includes/control_functions.php
${ZM_HTTPD_ROOT}/skins/mobile/includes/functions.php
${ZM_HTTPD_ROOT}/skins/mobile/includes/init.php
${ZM_HTTPD_ROOT}/skins/mobile/skin.php
${ZM_HTTPD_ROOT}/skins/mobile/views/console.php
${ZM_HTTPD_ROOT}/skins/mobile/views/css/console.css
${ZM_HTTPD_ROOT}/skins/mobile/views/devices.php
${ZM_HTTPD_ROOT}/skins/mobile/views/error.php
${ZM_HTTPD_ROOT}/skins/mobile/views/event.php
${ZM_HTTPD_ROOT}/skins/mobile/views/eventdetails.php
${ZM_HTTPD_ROOT}/skins/mobile/views/events.php
${ZM_HTTPD_ROOT}/skins/mobile/views/filter.php
${ZM_HTTPD_ROOT}/skins/mobile/views/frame.php
${ZM_HTTPD_ROOT}/skins/mobile/views/function.php
${ZM_HTTPD_ROOT}/skins/mobile/views/login.php
${ZM_HTTPD_ROOT}/skins/mobile/views/montage.php
${ZM_HTTPD_ROOT}/skins/mobile/views/state.php
${ZM_HTTPD_ROOT}/skins/mobile/views/video.php
${ZM_HTTPD_ROOT}/skins/mobile/views/watch.php
${ZM_HTTPD_ROOT}/skins/xml/includes/config.php
${ZM_HTTPD_ROOT}/skins/xml/includes/functions.php
${ZM_HTTPD_ROOT}/skins/xml/includes/init.php
${ZM_HTTPD_ROOT}/skins/xml/skin.php
${ZM_HTTPD_ROOT}/skins/xml/views/actions.php
${ZM_HTTPD_ROOT}/skins/xml/views/console.php
${ZM_HTTPD_ROOT}/skins/xml/views/none.php
${ZM_HTTPD_ROOT}/skins/xml/views/notfound.png
${ZM_HTTPD_ROOT}/temp
${ZM_HTTPD_ROOT}/tools/mootools/mootools-core-1.3.2-nc.js
${ZM_HTTPD_ROOT}/tools/mootools/mootools-core-1.3.2-yc.js
${ZM_HTTPD_ROOT}/tools/mootools/mootools-core.js
${ZM_HTTPD_ROOT}/tools/mootools/mootools-more-1.3.2.1-nc.js
${ZM_HTTPD_ROOT}/tools/mootools/mootools-more-1.3.2.1-yc.js
${ZM_HTTPD_ROOT}/tools/mootools/mootools-more.js
${ZM_HTTPD_ROOT}/views/file.php
${ZM_HTTPD_ROOT}/views/image.php
@pkgdir ${ZM_HTTPD_ROOT}/sounds
@pkgdir ${ZM_HTTPD_ROOT}/skins/mobile/lang
@pkgdir ${ZM_HTTPD_ROOT}/skins/mobile/graphics
@pkgdir ${ZM_HTTPD_ROOT}/skins/mobile/ajax
@pkgdir ${ZM_HTTPD_ROOT}/skins/classic/lang
@pkgdir ${ZM_HTTPD_ROOT}/skins/classic/ajax
