# mail/popa3d/PLIST
@comment $NetBSD: PLIST,v 1.9 2009/06/14 18:04:39 joerg Exp $
man/man8/popa3d.8
sbin/popa3d
share/doc/popa3d/DESIGN
share/doc/popa3d/INSTALL
share/doc/popa3d/LICENSE
share/examples/rc.d/popa3d
@comment empty dir to chroot to
@exec ${MKDIR} %D/share/empty
