# mail/popa3d/PLIST
@comment $NetBSD: PLIST,v 1.6 2003/12/07 19:18:17 wiz Exp $
etc/rc.d/popa3d
man/man8/popa3d.8
sbin/popa3d
share/doc/popa3d/DESIGN
share/doc/popa3d/INSTALL
share/doc/popa3d/LICENSE
@dirrm share/doc/popa3d
@comment empty dir to chroot to
@exec ${MKDIR} %D/share/empty
@dirrm share/empty
