# mail/popa3d/PLIST
@comment $NetBSD: PLIST,v 1.8 2005/05/02 20:34:00 reed Exp $
man/man8/popa3d.8
sbin/popa3d
share/doc/popa3d/DESIGN
share/doc/popa3d/INSTALL
share/doc/popa3d/LICENSE
share/examples/rc.d/popa3d
@dirrm share/doc/popa3d
@comment empty dir to chroot to
@exec ${MKDIR} %D/share/empty
@dirrm share/empty
