# net/miniupnpd/PLIST
@comment $NetBSD: PLIST,v 1.3 2009/06/14 18:09:34 joerg Exp $
sbin/miniupnpd
share/examples/miniupnpd/miniupnpd.conf
share/examples/rc.d/miniupnpd
