# www/wwwoffle/pkg/PLIST
@comment $NetBSD: PLIST,v 1.7 2000/11/20 05:12:40 itojun Exp $
bin/wwwoffle
bin/wwwoffle-ls
bin/wwwoffle-mv
bin/wwwoffle-read
bin/wwwoffle-rm
bin/wwwoffle-tools
bin/wwwoffle-write
etc/rc.d/wwwoffle.sh
etc/wwwoffle/wwwoffle.conf.default
@exec if [ ! -f %D/etc/wwwoffle/wwwoffle.conf ]; then cp -p %D/etc/wwwoffle/wwwoffle.conf.default %D/etc/wwwoffle/wwwoffle.conf; fi
libdata/wwwoffle-spool.tar.gz
man/man1/wwwoffle.1
man/man5/wwwoffle.conf.5
man/man8/wwwoffled.8
sbin/wwwoffled
share/doc/wwwoffle/CHANGES.CONF
share/doc/wwwoffle/CONVERT
share/doc/wwwoffle/COPYING
share/doc/wwwoffle/FAQ
share/doc/wwwoffle/INSTALL
share/doc/wwwoffle/NEWS
share/doc/wwwoffle/README
share/doc/wwwoffle/README.1st
share/doc/wwwoffle/README.CONF
share/doc/wwwoffle/README.PWD
share/doc/wwwoffle/README.URL
share/doc/wwwoffle/README.htdig
share/doc/wwwoffle/README.lang
share/doc/wwwoffle/README.udmsearch
share/doc/wwwoffle/README.win32
@dirrm share/doc/wwwoffle
@exec @GTAR@ xzCf / %D/libdata/wwwoffle-spool.tar.gz
@exec cd /var/wwwoffle/html && mv wwwoffle.pac wwwoffle.pac.dist && sed /PROXY/s/localhost/`hostname`/ <wwwoffle.pac.dist >wwwoffle.pac
@unexec cd /var/wwwoffle && rm -rf html* lasttime outgoing prevtime* temp
@unexec rmdir /var/wwwoffle/* 2>/dev/null && rmdir /var/wwwoffle ; :
@unexec if [ -d /var/wwwoffle ]; then echo "If you won't use the WWWOFFLE package anymore, remove the cache manually by"; echo '	rm -rf /var/wwwoffle'; fi
