# mail/exim/pkg/PLIST
@comment $NetBSD: PLIST,v 1.1 1999/09/07 13:21:03 ad Exp $
etc/exim/configure
etc/rc.d/exim.sh
sbin/exicyclog
sbin/exigrep
sbin/exim
sbin/exim_dbmbuild
sbin/exim_dumpdb
sbin/exim_fixdb
sbin/exim_lock
sbin/exim_tidydb
sbin/eximstats
sbin/exinext
sbin/exiqsumm
sbin/exiwhat
sbin/mailq
@unexec rmdir %D/etc/exim /var/log/exim 2>/dev/null || true
