# mail/exim/pkg/PLIST
@comment $NetBSD: PLIST,v 1.2 1999/09/18 20:28:27 ad Exp $
etc/exim/configure
etc/rc.d/exim.sh
etc/mailer.conf.exim
sbin/exicyclog
sbin/exigrep
sbin/exim
sbin/exim_dbmbuild
sbin/exim_dumpdb
sbin/exim_fixdb
sbin/exim_lock
sbin/exim_newaliases
sbin/exim_tidydb
sbin/eximstats
sbin/exinext
sbin/exiqsumm
sbin/exiwhat
@unexec rmdir %D/etc/exim /var/log/exim 2>/dev/null || true
