# mail/exim/pkg/PLIST
@comment $NetBSD: PLIST,v 1.6 1999/12/01 11:19:19 ad Exp $
etc/exim/configure
etc/exim/mailer.conf.exim
etc/rc.d/exim.sh
share/doc/exim/filter.html
share/doc/exim/filter_1.html
share/doc/exim/filter_10.html
share/doc/exim/filter_11.html
share/doc/exim/filter_12.html
share/doc/exim/filter_13.html
share/doc/exim/filter_14.html
share/doc/exim/filter_15.html
share/doc/exim/filter_16.html
share/doc/exim/filter_17.html
share/doc/exim/filter_18.html
share/doc/exim/filter_19.html
share/doc/exim/filter_2.html
share/doc/exim/filter_20.html
share/doc/exim/filter_21.html
share/doc/exim/filter_22.html
share/doc/exim/filter_23.html
share/doc/exim/filter_24.html
share/doc/exim/filter_25.html
share/doc/exim/filter_26.html
share/doc/exim/filter_27.html
share/doc/exim/filter_28.html
share/doc/exim/filter_29.html
share/doc/exim/filter_3.html
share/doc/exim/filter_30.html
share/doc/exim/filter_31.html
share/doc/exim/filter_32.html
share/doc/exim/filter_33.html
share/doc/exim/filter_34.html
share/doc/exim/filter_4.html
share/doc/exim/filter_5.html
share/doc/exim/filter_6.html
share/doc/exim/filter_7.html
share/doc/exim/filter_8.html
share/doc/exim/filter_9.html
share/doc/exim/filter_toc.html
share/doc/exim/oview.html
share/doc/exim/spec.html
share/doc/exim/spec_1.html
share/doc/exim/spec_10.html
share/doc/exim/spec_11.html
share/doc/exim/spec_12.html
share/doc/exim/spec_13.html
share/doc/exim/spec_14.html
share/doc/exim/spec_15.html
share/doc/exim/spec_16.html
share/doc/exim/spec_17.html
share/doc/exim/spec_18.html
share/doc/exim/spec_19.html
share/doc/exim/spec_2.html
share/doc/exim/spec_20.html
share/doc/exim/spec_21.html
share/doc/exim/spec_22.html
share/doc/exim/spec_23.html
share/doc/exim/spec_24.html
share/doc/exim/spec_25.html
share/doc/exim/spec_26.html
share/doc/exim/spec_27.html
share/doc/exim/spec_28.html
share/doc/exim/spec_29.html
share/doc/exim/spec_3.html
share/doc/exim/spec_30.html
share/doc/exim/spec_31.html
share/doc/exim/spec_32.html
share/doc/exim/spec_33.html
share/doc/exim/spec_34.html
share/doc/exim/spec_35.html
share/doc/exim/spec_36.html
share/doc/exim/spec_37.html
share/doc/exim/spec_38.html
share/doc/exim/spec_39.html
share/doc/exim/spec_4.html
share/doc/exim/spec_40.html
share/doc/exim/spec_41.html
share/doc/exim/spec_42.html
share/doc/exim/spec_43.html
share/doc/exim/spec_44.html
share/doc/exim/spec_45.html
share/doc/exim/spec_46.html
share/doc/exim/spec_47.html
share/doc/exim/spec_48.html
share/doc/exim/spec_49.html
share/doc/exim/spec_5.html
share/doc/exim/spec_50.html
share/doc/exim/spec_51.html
share/doc/exim/spec_52.html
share/doc/exim/spec_53.html
share/doc/exim/spec_54.html
share/doc/exim/spec_55.html
share/doc/exim/spec_56.html
share/doc/exim/spec_6.html
share/doc/exim/spec_7.html
share/doc/exim/spec_8.html
share/doc/exim/spec_9.html
share/doc/exim/spec_cindex.html
share/doc/exim/spec_findex.html
share/doc/exim/spec_toc.html
share/doc/exim/spectop.html
sbin/exicyclog
sbin/exigrep
sbin/exim
sbin/exim_dbmbuild
sbin/exim_dumpdb
sbin/exim_fixdb
sbin/exim_lock
sbin/exim_newaliases
sbin/exim_tidydb
sbin/eximstats
sbin/exinext
sbin/exiqsumm
sbin/exiwhat
man/man8/exim.8
@unexec rmdir %D/etc/exim %D/share/doc/exim /var/log/exim 2>/dev/null || true
