# net/rsync/PLIST
@comment $NetBSD: PLIST,v 1.5 2009/02/21 09:34:17 tron Exp $
bin/rsync
man/man1/rsync.1
man/man5/rsyncd.conf.5
share/doc/rsync/README
share/doc/rsync/tech_report.tex
share/examples/rc.d/rsyncd
@dirrm share/doc/rsync
