# www/varnish/PLIST
@comment $NetBSD: PLIST,v 1.7 2011/11/12 20:40:16 mspo Exp $
bin/varnishadm
bin/varnishhist
bin/varnishlog
bin/varnishncsa
bin/varnishreplay
bin/varnishsizes
bin/varnishstat
bin/varnishtest
bin/varnishtop
include/varnish/varnishapi.h
include/varnish/vcli.h
include/varnish/vsc.h
include/varnish/vsc_all.h
include/varnish/vsc_fields.h
include/varnish/vsl.h
include/varnish/vsl_tags.h
include/varnish/vsm.h
lib/libvarnishapi.la
lib/pkgconfig/varnishapi.pc
lib/varnish/libvarnish.la
lib/varnish/libvarnishcompat.la
lib/varnish/libvcl.la
lib/varnish/libvgz.la
lib/varnish/vmods/libvmod_std.la
man/man1/varnishadm.1
man/man1/varnishd.1
man/man1/varnishhist.1
man/man1/varnishlog.1
man/man1/varnishncsa.1
man/man1/varnishreplay.1
man/man1/varnishsizes.1
man/man1/varnishstat.1
man/man1/varnishtest.1
man/man1/varnishtop.1
man/man3/vmod_std.3
man/man7/varnish-cli.7
man/man7/vcl.7
man/man7/varnish-counters.7
sbin/varnishd
share/examples/rc.d/varnishd
share/examples/rc.d/varnishlog
share/examples/rc.d/varnishncsa
share/examples/varnish/default.vcl
