# www/squid3/PLIST
@comment $NetBSD: PLIST,v 1.1 2013/02/10 18:16:53 adam Exp $
bin/purge
bin/squidclient
${PLIST.ba_getpwnam}libexec/basic_getpwnam_auth
${PLIST.ba_MSNT}libexec/basic_msnt_auth
${PLIST.ba_NCSA}libexec/basic_ncsa_auth
${PLIST.ba_NIS}libexec/basic_nis_auth
${PLIST.ba_PAM}libexec/basic_pam_auth
libexec/cachemgr.cgi
${PLIST.da_ldap}libexec/digest_ldap_auth
${PLIST.diskd}libexec/diskd
${PLIST.eacl_unix_group}libexec/ext_unix_group_acl
libexec/helper-mux.pl
libexec/log_db_daemon
libexec/log_file_daemon
${PLIST.ta_kerberos}libexec/negotiate_kerberos_auth
${PLIST.ta_kerberos}libexec/negotiate_kerberos_auth_test
${PLIST.na_SMB}libexec/ntlm_auth
libexec/pinger
${PLIST.ba_LDAP}libexec/squid_ldap_auth
${PLIST.eacl_ldap_group}libexec/squid_ldap_group
${PLIST.unlinkd}libexec/unlinkd
libexec/url_fake_rewrite
libexec/url_fake_rewrite.sh
man/man1/squidclient.1
man/man8/cachemgr.cgi.8
${PLIST.ba_getpwnam}man/man8/basic_getpwnam_auth.8
${PLIST.ba_NCSA}man/man8/basic_ncsa_auth.8
${PLIST.ba_PAM}man/man8/basic_pam_auth.8
${PLIST.eacl_unix_group}man/man8/ext_unix_group_acl.8
${PLIST.ta_kerberos}man/man8/negotiate_kerberos_auth.8
man/man8/log_db_daemon.8
man/man8/squid.8
${PLIST.ba_LDAP}man/man8/squid_ldap_auth.8
${PLIST.eacl_ldap_group}man/man8/squid_ldap_group.8
sbin/squid
share/doc/squid/ChangeLog
share/doc/squid/README.FreeBSD
share/doc/squid/README.NetBSD
share/doc/squid/README.OpenBSD
share/doc/squid/README.Solaris
share/doc/squid/RELEASENOTES.html
share/doc/squid/debug-sections.txt
share/examples/rc.d/squid
share/examples/squid/cachemgr.conf
share/examples/squid/errorpage.css
share/examples/squid/mime.conf
share/examples/squid/msntauth.conf
share/examples/squid/squid.conf
share/squid/mib.txt
