# security/cvm/PLIST
@comment $NetBSD: PLIST,v 1.4 2009/06/14 18:13:27 joerg Exp $
bin/cvm-benchclient
bin/cvm-chain
bin/cvm-checkpassword
bin/cvm-pwfile
bin/cvm-qmail
bin/cvm-testclient
bin/cvm-unix
bin/cvm-v1benchclient
bin/cvm-v1checkpassword
bin/cvm-v1testclient
bin/cvm-vmailmgr
bin/cvm-vmailmgr-local
bin/cvm-vmailmgr-udp
include/cvm-sasl.h
include/cvm/client.h
include/cvm/credentials.h
include/cvm/errors.h
include/cvm/facts.h
include/cvm/module.h
include/cvm/protocol.h
include/cvm/sasl.h
include/cvm/v1client.h
include/cvm/v2client.h
lib/libcvm-client.a
lib/libcvm-client.la
lib/libcvm-client.${CVM_LIB_EXT}
lib/libcvm-command.la
lib/libcvm-local.la
lib/libcvm-module.la
lib/libcvm-qmail.la
lib/libcvm-sasl.la
lib/libcvm-sql.la
lib/libcvm-udp.la
lib/libcvm-v1client.la
lib/libcvm-v2client.la
share/examples/rc.d/cvm
