# graphics/nvtv/PLIST
@comment $NetBSD: PLIST,v 1.3 2005/05/02 20:33:59 reed Exp $
bin/nvtv
bin/nvtvd
share/doc/nvtv/FAQ
share/doc/nvtv/README
share/doc/nvtv/USAGE
share/doc/nvtv/chips.txt
share/doc/nvtv/overview.txt
share/doc/nvtv/timing.txt
man/man1/nvtv.1
man/man8/nvtvd.8
share/examples/rc.d/nvtvd
@dirrm share/doc/nvtv
