# chat/jabberd2/PLIST
@comment $NetBSD: PLIST,v 1.16 2013/12/12 23:01:02 schnoebe Exp $
bin/bdbdump.pl
bin/c2s
bin/jabberd
bin/pipe-auth.pl
bin/router
bin/s2s
bin/sm
lib/jabberd/authreg_anon.la
${PLIST.db}lib/jabberd/authreg_db.la
${PLIST.ldap}lib/jabberd/authreg_ldap.la
${PLIST.ldap}lib/jabberd/authreg_ldapfull.la
${PLIST.mysql}lib/jabberd/authreg_mysql.la
${PLIST.pam}lib/jabberd/authreg_pam.la
${PLIST.pgsql}lib/jabberd/authreg_pgsql.la
${PLIST.sqlite}lib/jabberd/authreg_sqlite.la
lib/jabberd/authreg_pipe.la
lib/jabberd/libstorage.la
lib/jabberd/mod_active.la
lib/jabberd/mod_amp.la
lib/jabberd/mod_announce.la
lib/jabberd/mod_deliver.la
lib/jabberd/mod_disco.la
lib/jabberd/mod_echo.la
lib/jabberd/mod_help.la
lib/jabberd/mod_iq-last.la
lib/jabberd/mod_iq-ping.la
lib/jabberd/mod_iq-private.la
lib/jabberd/mod_iq-time.la
lib/jabberd/mod_iq-vcard.la
lib/jabberd/mod_iq-version.la
lib/jabberd/mod_offline.la
lib/jabberd/mod_pep.la
lib/jabberd/mod_presence.la
lib/jabberd/mod_privacy.la
lib/jabberd/mod_roster-publish.la
lib/jabberd/mod_roster.la
lib/jabberd/mod_session.la
lib/jabberd/mod_status.la
lib/jabberd/mod_template-roster.la
lib/jabberd/mod_vacation.la
lib/jabberd/mod_validate.la
lib/jabberd/storage_fs.la
${PLIST.db}lib/jabberd/storage_db.la
${PLIST.ldap}lib/jabberd/storage_ldapvcard.la
${PLIST.mysql}lib/jabberd/storage_mysql.la
${PLIST.pgsql}lib/jabberd/storage_pgsql.la
${PLIST.sqlite}lib/jabberd/storage_sqlite.la
man/man8/c2s.8
man/man8/jabberd.8
man/man8/router.8
man/man8/s2s.8
man/man8/sm.8
share/doc/jabberd/ChangeLog
share/doc/jabberd/README
share/doc/jabberd/README.protocol
share/doc/jabberd/README.config
share/examples/jabberd/c2s.xml
share/examples/jabberd/db-jd14-2-jd2.sql
share/examples/jabberd/db-setup.mysql
share/examples/jabberd/db-setup.pgsql
share/examples/jabberd/db-setup.sqlite
share/examples/jabberd/db-update.mysql
share/examples/jabberd/db-update.sqlite
share/examples/jabberd/jabberd.cfg
share/examples/jabberd/roster.xml
share/examples/jabberd/router-filter.xml
share/examples/jabberd/router-users.xml
share/examples/jabberd/router.xml
share/examples/jabberd/s2s.xml
share/examples/jabberd/sm.xml
share/examples/rc.d/c2s
share/examples/rc.d/jabberd
share/examples/rc.d/router
share/examples/rc.d/s2s
share/examples/rc.d/sm
