# net/isc-dhcpd4/PLIST
@comment $NetBSD: PLIST,v 1.4 2013/12/13 16:34:17 jakllsch Exp $
bin/omshell
lib/libdhcpctl.a
lib/libomapi.a
man/man1/omshell.1
man/man3/dhcpctl.3
man/man3/omapi.3
man/man5/dhcpd.conf.5
man/man5/dhcpd.leases.5
man/man8/dhcpd.8
sbin/dhcpd
share/examples/isc-dhcpd/dhcpd.conf
share/examples/rc.d/isc_dhcpd
share/examples/rc.d/isc_dhcpd6
