# net/isc-dhcpd4/PLIST
@comment $NetBSD: PLIST,v 1.1.1.1 2008/02/13 22:06:15 adrianp Exp $
bin/omshell
lib/libdhcpctl.a
lib/libomapi.a
man/man1/omshell.1
man/man3/dhcpctl.3
man/man3/omapi.3
man/man5/dhcpd.conf.5
man/man5/dhcpd.leases.5
man/man8/dhcpd.8
sbin/dhcpd
share/examples/isc-dhcp-server/dhcpd.conf
share/examples/rc.d/isc_dhcpd
@dirrm share/examples/isc-dhcp-server
