# net/knot/PLIST
@comment $NetBSD: PLIST,v 1.4 2014/01/31 22:05:43 pettai Exp $
bin/kdig
bin/khost
bin/knsec3hash
bin/knsupdate
man/man1/kdig.1
man/man1/khost.1
man/man1/knsec3hash.1
man/man1/knsupdate.1
man/man5/knot.conf.5
man/man8/knotc.8
man/man8/knotd.8
sbin/knotc
sbin/knotd
share/examples/knot/example.com.zone
share/examples/knot/knot.sample.conf
share/examples/rc.d/knot
@pkgdir etc/knot
