# security/prelude-manager/PLIST
@comment $NetBSD: PLIST,v 1.6 2007/05/12 12:50:22 shannonjr Exp $
bin/prelude-manager
include/prelude-manager/prelude-manager.h
lib/prelude-manager/decodes/normalize.la
lib/prelude-manager/filters/idmef-criteria.la
lib/prelude-manager/filters/thresholding.a
lib/prelude-manager/filters/thresholding.la
lib/prelude-manager/filters/thresholding.so
lib/prelude-manager/reports/db.la
lib/prelude-manager/reports/debug.la
lib/prelude-manager/reports/relaying.la
lib/prelude-manager/reports/textmod.la
sbin/run-prelude-manager
share/examples/prelude-manager/prelude-manager.conf
share/examples/rc.d/preludemanager
@dirrm lib/prelude-manager/reports
@dirrm lib/prelude-manager/filters
@dirrm lib/prelude-manager/decodes
@dirrm lib/prelude-manager
@dirrm include/prelude-manager
