# www/wwwoffle/files/wwwoffled
#!/bin/sh
#
# $NetBSD: wwwoffled,v 1.2 2001/10/25 10:20:59 wiz Exp $
#

# PROVIDE: wwwoffled
# REQUIRE: DAEMON network

. /etc/rc.subr

name="wwwoffled"
command="@PREFIX@/sbin/${name}"
control_prog="@PREFIX@/bin/wwwoffle"
command_args="-c @PREFIX@/etc/wwwoffle/wwwoffle.conf"
required_files="@PREFIX@/etc/wwwoffle/wwwoffle.conf"
start_cmd=wwwoffled_start
stop_cmd="${control_prog} ${command_args} -kill"
autodial_cmd="${control_prog} -autodial ${command_args}"
offline_cmd="${control_prog} -offline ${command_args}"
online_cmd="${control_prog} -online ${command_args}"
extra_commands="autodial offline online"

wwwoffled_start() {
	eval "${command} ${command_args}"
	eval "${control_prog} -online ${command_args}"
	eval "${control_prog} -fetch ${command_args}"
}

load_rc_config $name
run_rc_command "$1"
