# misc/dnetc/PLIST
@comment $NetBSD: PLIST,v 1.4 2005/05/02 20:34:01 reed Exp $
bin/dnetc
man/man1/dnetc.1
share/doc/dnetc/CHANGES.txt
share/doc/dnetc/dnetc.txt
share/doc/dnetc/readme.1st
share/doc/dnetc/readme._ix
share/examples/rc.d/dnetc
@dirrm share/doc/dnetc
