# multimedia/gmediaserver/PLIST
@comment $NetBSD: PLIST,v 1.2 2013/02/16 15:10:46 tron Exp $
bin/gmediaserver
info/gmediaserver.info
man/man1/gmediaserver.1
share/examples/gmediaserver/gmediaserver.conf
share/examples/rc.d/gmediaserver
share/locale/sv/LC_MESSAGES/gmediaserver.mo
