# net/vsftpd/files/smf/manifest.xml
<?xml version='1.0'?>
<!DOCTYPE service_bundle SYSTEM '/usr/share/lib/xml/dtd/service_bundle.dtd.1'>
<service_bundle type='manifest' name='export'>
  <service name='@SMF_PREFIX@/@SMF_NAME@' type='service' version='1'>
    <dependency name='network' grouping='require_all' restart_on='refresh' type='service'>
      <service_fmri value='svc:/milestone/network:default'/>
    </dependency>
    <dependency name='filesystem' grouping='require_all' restart_on='refresh' type='service'>
      <service_fmri value='svc:/system/filesystem/local'/>
    </dependency>
    <instance name='default' enabled='false'>
       <dependency name='config' grouping='require_all' restart_on='none' type='path'>
	 <service_fmri value='file://localhost/@PKG_SYSCONFDIR@/vsftpd.conf' />
       </dependency>
       <exec_method type='method' name='start' exec='@PREFIX@/libexec/vsftpd %{vsftpd_params}' timeout_seconds='60'/>
       <exec_method type='method' name='stop'  exec=':kill' timeout_seconds='60'/>
       <property_group name='startd' type='framework'>
	 <propval name='duration' type='astring' value='child'/>
	 <propval name='ignore_error' type='astring' value='core,signal'/>
       </property_group>
       <property_group name='application' type='application'>
	 <propval name='vsftpd_params' type='astring' value='@PKG_SYSCONFDIR@/vsftpd.conf' />
       </property_group>
       <template>
	 <common_name>
	   <loctext xml:lang='C'>Very Secure FTP Daemon</loctext>
	 </common_name>
	 <documentation>
	   <manpage title='vsftpd' section='8' manpath='@PREFIX@/@PKGMANDIR@'/>
	 </documentation>
       </template>
    </instance>
  </service>
</service_bundle>
