# mail/prayer/PLIST
@comment $NetBSD: PLIST,v 1.7 2005/05/02 20:34:00 reed Exp $
sbin/prayer
sbin/prayer-cyclog
sbin/prayer-sem-prune
sbin/prayer-session
sbin/prayer-ssl-prune
share/examples/prayer/motd.html
share/examples/prayer/prayer-accountd.cf
share/examples/prayer/prayer.cf
share/examples/prayer/welcome.html
share/examples/rc.d/prayer
share/prayer/help/abook_add.html
share/prayer/help/abook_list.html
share/prayer/help/abook_search.html
share/prayer/help/abook_xfer.html
share/prayer/help/attachments.html
share/prayer/help/block.html
share/prayer/help/compose/large.html
share/prayer/help/compose/normal.html
share/prayer/help/compose/postponed.html
share/prayer/help/compose/role_list.html
share/prayer/help/copy.html
share/prayer/help/dictionary.html
share/prayer/help/display.html
share/prayer/help/favourites.html
share/prayer/help/filter.html
share/prayer/help/filter_select.html
share/prayer/help/folders.html
share/prayer/help/fullname.html
share/prayer/help/include.html
share/prayer/help/list.html
share/prayer/help/manage.html
share/prayer/help/passwd.html
share/prayer/help/prefs/compose1.html
share/prayer/help/prefs/compose2.html
share/prayer/help/prefs/display1.html
share/prayer/help/prefs/folder.html
share/prayer/help/prefs/general.html
share/prayer/help/quota.html
share/prayer/help/redirect.html
share/prayer/help/rename.html
share/prayer/help/reply.html
share/prayer/help/roles_entry.html
share/prayer/help/roles_list.html
share/prayer/help/search.html
share/prayer/help/sieve.html
share/prayer/help/spam.html
share/prayer/help/toolbar.html
share/prayer/help/transfer.html
share/prayer/help/upload_select.html
share/prayer/help/vacation.html
share/prayer/help/vaclog.html
share/prayer/help/welcome.html
share/prayer/icons/addressbook-inactive.gif
share/prayer/icons/addressbook.gif
share/prayer/icons/back.gif
share/prayer/icons/compose-inactive.gif
share/prayer/icons/compose.gif
share/prayer/icons/cross.gif
share/prayer/icons/dir.gif
share/prayer/icons/display-inactive.gif
share/prayer/icons/display.gif
share/prayer/icons/folders-inactive.gif
share/prayer/icons/folders.gif
share/prayer/icons/help-inactive.gif
share/prayer/icons/help.gif
share/prayer/icons/left-inactive.gif
share/prayer/icons/left.gif
share/prayer/icons/leftleft-inactive.gif
share/prayer/icons/leftleft.gif
share/prayer/icons/logout.gif
share/prayer/icons/mailbox-inactive.gif
share/prayer/icons/mailbox.gif
share/prayer/icons/manage-inactive.gif
share/prayer/icons/manage.gif
share/prayer/icons/msg-answered.gif
share/prayer/icons/msg-deleted.gif
share/prayer/icons/msg-seen.gif
share/prayer/icons/msg-unseen.gif
share/prayer/icons/off.gif
share/prayer/icons/on.gif
share/prayer/icons/paperclip-small.gif
share/prayer/icons/right-inactive.gif
share/prayer/icons/right.gif
share/prayer/icons/rightright-inactive.gif
share/prayer/icons/rightright.gif
share/prayer/icons/tick.gif
share/prayer/icons/univheader.gif
@dirrm share/prayer/icons
@dirrm share/prayer/help/prefs
@dirrm share/prayer/help/compose
@dirrm share/prayer/help
@dirrm share/prayer
