# mail/prayer/PLIST
@comment $NetBSD: PLIST,v 1.11 2010/06/08 12:34:38 adam Exp $
man/man5/prayer.cf.5
man/man8/prayer-session.8
man/man8/prayer-ssl-prune.8
man/man8/prayer.8
sbin/prayer
sbin/prayer-chroot
sbin/prayer-cyclog
sbin/prayer-db-prune
sbin/prayer-sem-prune
sbin/prayer-session
sbin/prayer-ssl-prune
share/examples/prayer/motd.html
share/examples/prayer/prayer-accountd.cf
share/examples/prayer/prayer.cf
share/examples/rc.d/prayer
share/prayer/icons/addressbook-inactive.gif
share/prayer/icons/addressbook.gif
share/prayer/icons/back.gif
share/prayer/icons/bg-home-menu-top.gif
share/prayer/icons/bg-nav-primary.gif
share/prayer/icons/bg-tab-left-active.gif
share/prayer/icons/bg-tab-left-over.gif
share/prayer/icons/bg-tab-left.gif
share/prayer/icons/bg-tab-right-active.gif
share/prayer/icons/bg-tab-right-over.gif
share/prayer/icons/bg-tab-right.gif
share/prayer/icons/blank.gif
share/prayer/icons/button-search.gif
share/prayer/icons/button2-bm.png
share/prayer/icons/compose-inactive.gif
share/prayer/icons/compose.gif
share/prayer/icons/cross.gif
share/prayer/icons/dir.gif
share/prayer/icons/display-inactive.gif
share/prayer/icons/display.gif
share/prayer/icons/divider-site-tools-small.gif
share/prayer/icons/down.gif
share/prayer/icons/folders-inactive.gif
share/prayer/icons/folders.gif
share/prayer/icons/help-inactive.gif
share/prayer/icons/help.gif
share/prayer/icons/icon-alert-normal.gif
share/prayer/icons/icon-alert-warning.gif
share/prayer/icons/icon-arrow-home-menu-orange.gif
share/prayer/icons/icon-arrow-home-menu.gif
share/prayer/icons/icon-arrow-nav-breadcrumb.gif
share/prayer/icons/icon-arrow-nav-primary-down.gif
share/prayer/icons/icon-arrow-nav-primary-up.gif
share/prayer/icons/icon-arrow-nav-primary.gif
share/prayer/icons/icon-home.gif
share/prayer/icons/icon-quick-links.gif
share/prayer/icons/icon-section.gif
share/prayer/icons/identifier.gif
share/prayer/icons/identifier800.gif
share/prayer/icons/ifmlogo-small.gif
share/prayer/icons/ifmlogo.gif
share/prayer/icons/left-inactive.gif
share/prayer/icons/left.gif
share/prayer/icons/leftleft-inactive.gif
share/prayer/icons/leftleft.gif
share/prayer/icons/logout.gif
share/prayer/icons/mailbox-inactive.gif
share/prayer/icons/mailbox.gif
share/prayer/icons/manage-inactive.gif
share/prayer/icons/manage.gif
share/prayer/icons/msg-answered.gif
share/prayer/icons/msg-deleted.gif
share/prayer/icons/msg-seen.gif
share/prayer/icons/msg-seen.png
share/prayer/icons/msg-unseen.gif
share/prayer/icons/next-white.gif
share/prayer/icons/next.gif
share/prayer/icons/off.gif
share/prayer/icons/on.gif
share/prayer/icons/paperclip-small.gif
share/prayer/icons/prev-white.gif
share/prayer/icons/previous.gif
share/prayer/icons/raven.gif
share/prayer/icons/right-inactive.gif
share/prayer/icons/right.gif
share/prayer/icons/rightright-inactive.gif
share/prayer/icons/rightright.gif
share/prayer/icons/rssfeed.gif
share/prayer/icons/sort-d.gif
share/prayer/icons/sort-u.gif
share/prayer/icons/srssfeed.gif
share/prayer/icons/tick.gif
share/prayer/icons/ucs-id-small.jpg
share/prayer/icons/ucs-s.gif
share/prayer/icons/ucs.gif
share/prayer/icons/univheader.gif
share/prayer/static/blue.css
share/prayer/static/camstyle.css
share/prayer/static/common.css
share/prayer/static/content.css
share/prayer/static/forms.css
share/prayer/static/gray.css
share/prayer/static/green.css
share/prayer/static/help.css
share/prayer/static/high.css
share/prayer/static/ie6.css
share/prayer/static/ie7.css
share/prayer/static/layout.css
share/prayer/static/links.css
share/prayer/static/minmax.js
share/prayer/static/print.css
share/prayer/static/reset.css
share/prayer/static/rssfeed.css
share/prayer/static/typography.css
share/prayer/static/yellow.css
