# net/isc-dhcpd4/files/smf/manifest.xml
<?xml version='1.0'?>
<!DOCTYPE service_bundle SYSTEM '/usr/share/lib/xml/dtd/service_bundle.dtd.1'>
<service_bundle type='manifest' name='export'>
  <service name='@SMF_PREFIX@/@SMF_NAME@' type='service' version='1'>
    <create_default_instance enabled='false' />
    <single_instance />
    <dependency name='network' grouping='require_all' restart_on='error' type='service'>
      <service_fmri value='svc:/milestone/network:default' />
    </dependency>
    <dependency name='filesystem-local' grouping='require_all' restart_on='none' type='service'>
      <service_fmri value='svc:/system/filesystem/local:default' />
    </dependency>
    <dependency name='config-file' grouping='require_all' restart_on='refresh' type='path'>
      <service_fmri value='file://@PKG_SYSCONFDIR@/dhcpd.conf'/>
    </dependency>
    <method_context>
      <method_credential user='root' group='root' />
    </method_context>
    <exec_method type='method' name='start' exec='@PREFIX@/@SMF_METHOD_FILE.isc-dhcpd@ start' timeout_seconds='30' />
    <exec_method type='method' name='stop' exec='@PREFIX@/@SMF_METHOD_FILE.isc-dhcpd@ stop' timeout_seconds='30' />
    <exec_method type='method' name='refresh' exec='@PREFIX@/@SMF_METHOD_FILE.isc-dhcpd@ refresh' timeout_seconds='30' />
    <template>
      <common_name>
        <loctext xml:lang='C'>ISC DHCP Server</loctext>
      </common_name>
    </template>
  </service>
</service_bundle>
