# parallel/ganglia-monitor-core/PLIST
@comment $NetBSD: PLIST,v 1.3 2005/05/02 20:34:04 reed Exp $
bin/ganglia-config
bin/gmetric
bin/gstat
include/ganglia.h
lib/libganglia.la
man/man1/gmetad.1
man/man1/gmetric.1
man/man1/gmond.1
man/man1/gstat.1
sbin/gmetad
sbin/gmond
share/examples/ganglia/gmetad.conf
share/examples/ganglia/gmond.conf
share/examples/rc.d/gmetad
share/examples/rc.d/gmond
@dirrm share/examples/ganglia
