# chat/unrealircd/PLIST
@comment $NetBSD: PLIST,v 1.1.1.1 2005/11/29 19:17:03 adrianp Exp $
sbin/ircd
share/doc/html/unrealircd/unreal32docs.html
share/doc/unrealircd/Authors
share/doc/unrealircd/LICENSE
share/doc/unrealircd/coding-guidelines
share/doc/unrealircd/tao.of.irc
share/examples/rc.d/unrealircd
share/examples/unrealircd/example.conf
share/unrealircd/aliases/aliases.conf
share/unrealircd/aliases/anope.conf
share/unrealircd/aliases/auspice.conf
share/unrealircd/aliases/cygnus.conf
share/unrealircd/aliases/epona.conf
share/unrealircd/aliases/generic.conf
share/unrealircd/aliases/genericstats.conf
share/unrealircd/aliases/ircservices.conf
share/unrealircd/aliases/operstats.conf
share/unrealircd/badwords.channel.conf
share/unrealircd/badwords.message.conf
share/unrealircd/badwords.quit.conf
share/unrealircd/help.conf
share/unrealircd/modules/cloak.so
share/unrealircd/modules/commands.so
share/unrealircd/modules/m_addline.so
share/unrealircd/modules/m_addomotd.so
share/unrealircd/modules/m_admin.so
share/unrealircd/modules/m_adminchat.so
share/unrealircd/modules/m_akill.so
share/unrealircd/modules/m_away.so
share/unrealircd/modules/m_botmotd.so
share/unrealircd/modules/m_chatops.so
share/unrealircd/modules/m_chghost.so
share/unrealircd/modules/m_chgident.so
share/unrealircd/modules/m_chgname.so
share/unrealircd/modules/m_close.so
share/unrealircd/modules/m_connect.so
share/unrealircd/modules/m_cycle.so
share/unrealircd/modules/m_dccallow.so
share/unrealircd/modules/m_dccdeny.so
share/unrealircd/modules/m_eos.so
share/unrealircd/modules/m_globops.so
share/unrealircd/modules/m_guest.so
share/unrealircd/modules/m_help.so
share/unrealircd/modules/m_htm.so
share/unrealircd/modules/m_invite.so
share/unrealircd/modules/m_ison.so
share/unrealircd/modules/m_join.so
share/unrealircd/modules/m_kick.so
share/unrealircd/modules/m_kill.so
share/unrealircd/modules/m_knock.so
share/unrealircd/modules/m_lag.so
share/unrealircd/modules/m_links.so
share/unrealircd/modules/m_list.so
share/unrealircd/modules/m_locops.so
share/unrealircd/modules/m_lusers.so
share/unrealircd/modules/m_map.so
share/unrealircd/modules/m_message.so
share/unrealircd/modules/m_mkpasswd.so
share/unrealircd/modules/m_mode.so
share/unrealircd/modules/m_motd.so
share/unrealircd/modules/m_nachat.so
share/unrealircd/modules/m_netinfo.so
share/unrealircd/modules/m_nick.so
share/unrealircd/modules/m_oper.so
share/unrealircd/modules/m_opermotd.so
share/unrealircd/modules/m_part.so
share/unrealircd/modules/m_pass.so
share/unrealircd/modules/m_pingpong.so
share/unrealircd/modules/m_protoctl.so
share/unrealircd/modules/m_quit.so
share/unrealircd/modules/m_rakill.so
share/unrealircd/modules/m_rping.so
share/unrealircd/modules/m_rules.so
share/unrealircd/modules/m_sajoin.so
share/unrealircd/modules/m_samode.so
share/unrealircd/modules/m_sapart.so
share/unrealircd/modules/m_sdesc.so
share/unrealircd/modules/m_sendsno.so
share/unrealircd/modules/m_sendumode.so
share/unrealircd/modules/m_server.so
share/unrealircd/modules/m_sethost.so
share/unrealircd/modules/m_setident.so
share/unrealircd/modules/m_setname.so
share/unrealircd/modules/m_silence.so
share/unrealircd/modules/m_sjoin.so
share/unrealircd/modules/m_sqline.so
share/unrealircd/modules/m_squit.so
share/unrealircd/modules/m_stats.so
share/unrealircd/modules/m_svsfline.so
share/unrealircd/modules/m_svsjoin.so
share/unrealircd/modules/m_svskill.so
share/unrealircd/modules/m_svsmode.so
share/unrealircd/modules/m_svsmotd.so
share/unrealircd/modules/m_svsnick.so
share/unrealircd/modules/m_svsnline.so
share/unrealircd/modules/m_svsnoop.so
share/unrealircd/modules/m_svso.so
share/unrealircd/modules/m_svspart.so
share/unrealircd/modules/m_svssilence.so
share/unrealircd/modules/m_svssno.so
share/unrealircd/modules/m_svswatch.so
share/unrealircd/modules/m_swhois.so
share/unrealircd/modules/m_time.so
share/unrealircd/modules/m_tkl.so
share/unrealircd/modules/m_topic.so
share/unrealircd/modules/m_trace.so
share/unrealircd/modules/m_tsctl.so
share/unrealircd/modules/m_umode2.so
share/unrealircd/modules/m_undccdeny.so
share/unrealircd/modules/m_unkline.so
share/unrealircd/modules/m_user.so
share/unrealircd/modules/m_userip.so
share/unrealircd/modules/m_unsqline.so
share/unrealircd/modules/m_unzline.so
share/unrealircd/modules/m_userhost.so
share/unrealircd/modules/m_vhost.so
share/unrealircd/modules/m_wallops.so
share/unrealircd/modules/m_watch.so
share/unrealircd/modules/m_who.so
share/unrealircd/modules/m_whois.so
share/unrealircd/modules/m_whowas.so
share/unrealircd/networks/awesomechristians.network
share/unrealircd/networks/axenet.network
share/unrealircd/networks/bunker7.network
share/unrealircd/networks/burnnet.network
share/unrealircd/networks/cabonet.network
share/unrealircd/networks/chatcrap.network
share/unrealircd/networks/chatuniverse.network
share/unrealircd/networks/ctcp.network
share/unrealircd/networks/darkkaos.network
share/unrealircd/networks/digitalirc.network
share/unrealircd/networks/discussioni.network
share/unrealircd/networks/dragonwings.network
share/unrealircd/networks/gamescafe.network
share/unrealircd/networks/german-elite.network
share/unrealircd/networks/german-global-irc.network
share/unrealircd/networks/global-irc.network
share/unrealircd/networks/globalchat.network
share/unrealircd/networks/icechat.network
share/unrealircd/networks/infinity.network
share/unrealircd/networks/ircsystems.network
share/unrealircd/networks/isno.network
share/unrealircd/networks/l33t-irc.network
share/unrealircd/networks/lcirc.network
share/unrealircd/networks/makenet
share/unrealircd/networks/networks.ndx
share/unrealircd/networks/outsiderz.network
share/unrealircd/networks/phazenet.network
share/unrealircd/networks/stormdancing.network
share/unrealircd/networks/template.network
share/unrealircd/networks/thainet.network
share/unrealircd/networks/unitedirc-org.network
share/unrealircd/networks/unreal-test.network
share/unrealircd/networks/wazzza.network
share/unrealircd/networks/x-irc.network
share/unrealircd/networks/zirc.network
share/unrealircd/spamfilter.conf
@dirrm share/unrealircd/networks
@dirrm share/unrealircd/modules
@dirrm share/unrealircd/aliases
@dirrm share/unrealircd
@dirrm share/examples/unrealircd
@dirrm share/doc/unrealircd
@dirrm share/doc/html/unrealircd
