# www/hiawatha/PLIST
@comment $NetBSD: PLIST,v 1.3 2005/05/02 20:34:07 reed Exp $
lib/java/hiawatha/hiawatha.jar
lib/java/hiawatha/xsls.jar
share/examples/hiawatha/simple.cfg
share/examples/hiawatha/pault_com.cfg
share/examples/hiawatha/hiawatha.cfg
share/examples/hiawatha/pault.com/css/home-page.css
share/examples/hiawatha/pault.com/css/resume.css
share/examples/hiawatha/pault.com/css/xslscript.css
share/examples/hiawatha/pault.com/images/connector.gif
share/examples/hiawatha/pault.com/images/dtdgen1.gif
share/examples/hiawatha/pault.com/images/dtdgen2.gif
share/examples/hiawatha/pault.com/images/fish.gif
share/examples/hiawatha/pault.com/images/hiawatha.gif
share/examples/hiawatha/pault.com/images/hiawatha.jpeg
share/examples/hiawatha/pault.com/images/kay.gif
share/examples/hiawatha/pault.com/images/pagebg.gif
share/examples/hiawatha/pault.com/images/paultlogo.gif
share/examples/hiawatha/pault.com/images/pxsl.gif
share/examples/hiawatha/pault.com/images/resume.jpg
share/examples/hiawatha/pault.com/images/some2xml.gif
share/examples/hiawatha/pault.com/images/ux.gif
share/examples/hiawatha/pault.com/images/w3cval.gif
share/examples/hiawatha/pault.com/index-plain.xsls
share/examples/hiawatha/pault.com/index.xml
share/examples/hiawatha/pault.com/index.xsls
share/examples/hiawatha/pault.com/prod/Hiawatha/index-plain.xsls
share/examples/hiawatha/pault.com/prod/Hiawatha/index.xml
share/examples/hiawatha/pault.com/prod/Hiawatha/index.xsls
share/examples/hiawatha/pault.com/prod/Hiawatha/common.xsls
share/examples/hiawatha/pault.com/prod/Ux/index-plain.xsls
share/examples/hiawatha/pault.com/prod/Ux/index.xml
share/examples/hiawatha/pault.com/prod/Ux/index.xsls
share/examples/hiawatha/pault.com/prod/Ux/uxdoc/aliases.xml
share/examples/hiawatha/pault.com/prod/Ux/uxdoc/commands.xml
share/examples/hiawatha/pault.com/prod/Ux/uxdoc/concepts.xml
share/examples/hiawatha/pault.com/prod/Ux/uxdoc/history.xml
share/examples/hiawatha/pault.com/prod/Ux/uxdoc/index.xml
share/examples/hiawatha/pault.com/prod/Ux/uxdoc/index.xsls
share/examples/hiawatha/pault.com/prod/Ux/uxdoc/kernel.xml
share/examples/hiawatha/pault.com/prod/Ux/uxdoc/links.xsls
share/examples/hiawatha/pault.com/prod/Ux/uxdoc/master.xsl
share/examples/hiawatha/pault.com/prod/Ux/uxdoc/ux-beans.xml
share/examples/hiawatha/pault.com/prod/XSLScript/common.xsls
share/examples/hiawatha/pault.com/prod/XSLScript/index-plain.xsls
share/examples/hiawatha/pault.com/prod/XSLScript/index.xml
share/examples/hiawatha/pault.com/prod/XSLScript/index.xsls
share/examples/hiawatha/pault.com/resume-plain.xsls
share/examples/hiawatha/pault.com/resume.xml
share/examples/hiawatha/pault.com/resume.xsls
share/examples/hiawatha/pault.com/unsup/Phttpd/index-plain.xsls
share/examples/hiawatha/pault.com/unsup/Phttpd/index.xml
share/examples/hiawatha/pault.com/unsup/Phttpd/index.xsls
share/examples/hiawatha/pault.com/unsup/PJC/index-plain.xsls
share/examples/hiawatha/pault.com/unsup/PJC/index.xml
share/examples/hiawatha/pault.com/unsup/PJC/index.xsls
share/examples/hiawatha/pault.com/unsup/Pxsl/index-plain.xsls
share/examples/hiawatha/pault.com/unsup/Pxsl/index.xml
share/examples/hiawatha/pault.com/unsup/Pxsl/index.xsls
share/examples/hiawatha/pault.com/unsup/Some2XML/index-plain.xsls
share/examples/hiawatha/pault.com/unsup/Some2XML/index.xml
share/examples/hiawatha/pault.com/unsup/Some2XML/index.xsls
share/examples/hiawatha/pault.com/xsl-bin/cat.xsls
share/examples/hiawatha/pault.com/xsl-include/dbg.xsls
share/examples/hiawatha/pault.com/xsl-include/font.xsls
share/examples/hiawatha/pault.com/xsl-include/htmlpage.xsls
share/examples/hiawatha/pault.com/xsl-include/link.xsls
share/examples/hiawatha/pault.com/xsl-include/page-data.xsls
share/examples/hiawatha/pault.com/xsl-include/page-footer.xsls
share/examples/hiawatha/pault.com/xsl-include/page-logo.xsls
share/examples/hiawatha/pault.com/xsl-include/page-menu.xsls
share/examples/hiawatha/pault.com/xsl-include/page.xsls
share/examples/hiawatha/pault.com/xsl-include/plain-htmlpage.xsls
share/examples/hiawatha/pault.com/xsl-include/render-product.xsl
share/examples/hiawatha/pault.com/xsl-include/table.xsls
share/examples/hiawatha/pault.com/xsl-include/widget.xsls
share/examples/hiawatha/WEBROOT/index.xml
share/examples/hiawatha/WEBROOT/index.xsl
share/examples/hiawatha/WEBROOT/images/sitedrivenby.gif
share/examples/rc.d/hiawatha
@dirrm share/examples/hiawatha/pault.com/xsl-include
@dirrm share/examples/hiawatha/pault.com/xsl-bin
@dirrm share/examples/hiawatha/pault.com/unsup/Some2XML
@dirrm share/examples/hiawatha/pault.com/unsup/Pxsl
@dirrm share/examples/hiawatha/pault.com/unsup/Phttpd
@dirrm share/examples/hiawatha/pault.com/unsup/PJC
@dirrm share/examples/hiawatha/pault.com/unsup
@dirrm share/examples/hiawatha/pault.com/prod/XSLScript
@dirrm share/examples/hiawatha/pault.com/prod/Ux/uxdoc
@dirrm share/examples/hiawatha/pault.com/prod/Ux
@dirrm share/examples/hiawatha/pault.com/prod/Hiawatha
@dirrm share/examples/hiawatha/pault.com/prod
@dirrm share/examples/hiawatha/pault.com/images
@dirrm share/examples/hiawatha/pault.com/css
@dirrm share/examples/hiawatha/pault.com
@dirrm share/examples/hiawatha/WEBROOT/images
@dirrm share/examples/hiawatha/WEBROOT
@dirrm share/examples/hiawatha
@dirrm lib/java/hiawatha
