# net/net-snmp/files/smf/manifest.xml
<?xml version='1.0'?>
<!DOCTYPE service_bundle SYSTEM '/usr/share/lib/xml/dtd/service_bundle.dtd.1'>
<service_bundle type='manifest' name='export'>
  <service name='@SMF_PREFIX@/@SMF_NAME@' type='service' version='1'>
    <dependency name='milestone' grouping='require_all' restart_on='none' type='service'>
      <service_fmri value='svc:/milestone/sysconfig'/>
    </dependency>
    <dependency name='fs-local' grouping='require_all' restart_on='none' type='service'>
      <service_fmri value='svc:/system/filesystem/local'/>
    </dependency>
    <dependency name='name-services' grouping='optional_all' restart_on='none' type='service'>
      <service_fmri value='svc:/milestone/name-services'/>
    </dependency>
    <dependency name='system-log' grouping='optional_all' restart_on='none' type='service'>
      <service_fmri value='svc:/system/system-log'/>
    </dependency>
    <dependency name='rstat' grouping='optional_all' restart_on='none' type='service'>
      <service_fmri value='svc:/network/rpc/rstat'/>
    </dependency>
    <dependency name='cryptosvc' grouping='require_all' restart_on='restart' type='service'>
      <service_fmri value='svc:/system/cryptosvc'/>
    </dependency>
    <dependency name='network' grouping='require_all' restart_on='restart' type='service'>
      <service_fmri value='svc:/milestone/network'/>
    </dependency>
    <instance name='default' enabled='false'>
      <dependency name='config-file' grouping='require_all' restart_on='refresh' type='path'>
        <service_fmri value='file://localhost/@PKG_SYSCONFDIR@/snmpd.conf'/>
      </dependency>
      <exec_method name='start' type='method' exec='@PREFIX@/sbin/snmpd -c %{config_file} -Ls d' timeout_seconds='60'/>
      <exec_method name='stop' type='method' exec=':kill' timeout_seconds='60'/>
      <property_group name="application" type="application">
        <propval name="config_file" type="astring" value="@PKG_SYSCONFDIR@/snmpd.conf" />
      </property_group>
      <template>
        <common_name>
          <loctext xml:lang='C'>net-snmp SNMP daemon</loctext>
        </common_name>
        <documentation>
          <manpage title='snmpd' section='1M' manpath='@PREFIX@/@PKGMANDIR@/'/>
        </documentation>
      </template>
    </instance>
    <instance name='snmptrapd' enabled='false'>
      <dependency name='config-file' grouping='require_all' restart_on='refresh' type='path'>
        <service_fmri value='file://localhost/@PKG_SYSCONFDIR@/snmptrapd.conf'/>
      </dependency>
      <exec_method name='start' type='method' exec='@PREFIX@/sbin/snmptrapd -c %{config_file} -Ls 1' timeout_seconds='60'/>
      <exec_method name='stop' type='method' exec=':kill' timeout_seconds='60'/>
      <property_group name="application" type="application">
        <propval name="config_file" type="astring" value="@PKG_SYSCONFDIR@/snmptrapd.conf" />
      </property_group>
      <template>
        <common_name>
          <loctext xml:lang='C'>net-snmp SNMP trapper daemon</loctext>
        </common_name>
        <documentation>
          <manpage title='snmpd' section='1M' manpath='@PREFIX@/@PKGMANDIR@/'/>
        </documentation>
      </template>
    </instance>
  </service>
</service_bundle>
