# mail/spamassassin/PLIST
@comment $NetBSD: PLIST,v 1.11 2003/11/03 00:56:26 heinz Exp $
@comment Most of this PLIST is auto-generated; see Makefile.
etc/rc.d/spamd
share/doc/html/spamassassin/ninjabutton.png
@comment @dirrm share/doc/html/spamassassin
share/doc/spamassassin/2.60_change_summary
share/examples/spamassassin/local.cf
share/examples/spamassassin/netbsd_lists.cf
share/examples/spamassassin/procmailrc.example
share/spamassassin/languages
share/spamassassin/triplets.txt
share/spamassassin/user_prefs.template
@comment @dirrm share/examples/spamassassin
