# mail/spamassassin/PLIST
@comment $NetBSD: PLIST,v 1.16 2005/05/02 20:34:01 reed Exp $
@comment Most of this PLIST is auto-generated; see Makefile.
share/examples/rc.d/spamd
share/examples/spamassassin/init.pre
share/examples/spamassassin/local.cf
share/examples/spamassassin/netbsd_lists.cf
share/examples/spamassassin/procmailrc.example
share/spamassassin/languages
share/spamassassin/triplets.txt
share/spamassassin/user_prefs.template
