# mail/spamassassin/PLIST
@comment $NetBSD: PLIST,v 1.12 2003/12/27 19:16:04 heinz Exp $
@comment Most of this PLIST is auto-generated; see Makefile.
etc/rc.d/spamd
share/doc/html/spamassassin/ninjabutton.png
@comment @dirrm share/doc/html/spamassassin
share/doc/spamassassin/change_summary.txt
share/examples/spamassassin/local.cf
share/examples/spamassassin/netbsd_lists.cf
share/examples/spamassassin/procmailrc.example
share/spamassassin/languages
share/spamassassin/triplets.txt
share/spamassassin/user_prefs.template
@comment @dirrm share/examples/spamassassin
