# net/nagios-base/PLIST
@comment $NetBSD: PLIST,v 1.13 2013/08/05 10:53:34 obache Exp $
bin/nagios
bin/nagiostats
bin/p1.pl
libexec/nagios/cgi-bin/avail.cgi
libexec/nagios/cgi-bin/cmd.cgi
libexec/nagios/cgi-bin/config.cgi
libexec/nagios/cgi-bin/extinfo.cgi
libexec/nagios/cgi-bin/histogram.cgi
libexec/nagios/cgi-bin/history.cgi
libexec/nagios/cgi-bin/notifications.cgi
libexec/nagios/cgi-bin/outages.cgi
libexec/nagios/cgi-bin/showlog.cgi
libexec/nagios/cgi-bin/status.cgi
libexec/nagios/cgi-bin/statusmap.cgi
libexec/nagios/cgi-bin/statuswml.cgi
libexec/nagios/cgi-bin/statuswrl.cgi
libexec/nagios/cgi-bin/summary.cgi
libexec/nagios/cgi-bin/tac.cgi
libexec/nagios/cgi-bin/trends.cgi
share/examples/nagios/cgi.cfg-sample
share/examples/nagios/nagios.cfg-sample
share/examples/nagios/nagios.conf
share/examples/nagios/objects/commands.cfg-sample
share/examples/nagios/objects/contacts.cfg-sample
share/examples/nagios/objects/localhost.cfg-sample
share/examples/nagios/objects/printer.cfg-sample
share/examples/nagios/objects/switch.cfg-sample
share/examples/nagios/objects/templates.cfg-sample
share/examples/nagios/objects/timeperiods.cfg-sample
share/examples/nagios/objects/windows.cfg-sample
share/examples/nagios/resource.cfg-sample
share/examples/rc.d/nagios
share/nagios/config.inc.php
share/nagios/contexthelp/A1.html
share/nagios/contexthelp/A2.html
share/nagios/contexthelp/A3.html
share/nagios/contexthelp/A4.html
share/nagios/contexthelp/A5.html
share/nagios/contexthelp/A6.html
share/nagios/contexthelp/A7.html
share/nagios/contexthelp/B1.html
share/nagios/contexthelp/C1.html
share/nagios/contexthelp/D1.html
share/nagios/contexthelp/E1.html
share/nagios/contexthelp/F1.html
share/nagios/contexthelp/G1.html
share/nagios/contexthelp/G2.html
share/nagios/contexthelp/G3.html
share/nagios/contexthelp/G4.html
share/nagios/contexthelp/G5.html
share/nagios/contexthelp/G6.html
share/nagios/contexthelp/H1.html
share/nagios/contexthelp/H2.html
share/nagios/contexthelp/H3.html
share/nagios/contexthelp/H4.html
share/nagios/contexthelp/H5.html
share/nagios/contexthelp/H6.html
share/nagios/contexthelp/H7.html
share/nagios/contexthelp/H8.html
share/nagios/contexthelp/I1.html
share/nagios/contexthelp/I2.html
share/nagios/contexthelp/I3.html
share/nagios/contexthelp/I4.html
share/nagios/contexthelp/I5.html
share/nagios/contexthelp/I6.html
share/nagios/contexthelp/I7.html
share/nagios/contexthelp/I8.html
share/nagios/contexthelp/I9.html
share/nagios/contexthelp/J1.html
share/nagios/contexthelp/K1.html
share/nagios/contexthelp/L1.html
share/nagios/contexthelp/L10.html
share/nagios/contexthelp/L11.html
share/nagios/contexthelp/L12.html
share/nagios/contexthelp/L13.html
share/nagios/contexthelp/L2.html
share/nagios/contexthelp/L3.html
share/nagios/contexthelp/L4.html
share/nagios/contexthelp/L5.html
share/nagios/contexthelp/L6.html
share/nagios/contexthelp/L7.html
share/nagios/contexthelp/L8.html
share/nagios/contexthelp/L9.html
share/nagios/contexthelp/M1.html
share/nagios/contexthelp/M2.html
share/nagios/contexthelp/M3.html
share/nagios/contexthelp/M4.html
share/nagios/contexthelp/M5.html
share/nagios/contexthelp/M6.html
share/nagios/contexthelp/N1.html
share/nagios/contexthelp/N2.html
share/nagios/contexthelp/N3.html
share/nagios/contexthelp/N4.html
share/nagios/contexthelp/N5.html
share/nagios/contexthelp/N6.html
share/nagios/contexthelp/N7.html
share/nagios/images/Nagios-clearbg.png
share/nagios/images/NagiosEnterprises-whitebg-112x46.png
share/nagios/images/ack.gif
share/nagios/images/action-graph.gif
share/nagios/images/action-nagios.gif
share/nagios/images/action-orig.gif
share/nagios/images/action.gif
share/nagios/images/b_first2.png
share/nagios/images/b_last2.png
share/nagios/images/b_next2.png
share/nagios/images/b_prev2.png
share/nagios/images/command.png
share/nagios/images/comment.gif
share/nagios/images/contexthelp1.gif
share/nagios/images/contexthelp2.gif
share/nagios/images/critical.png
share/nagios/images/delay.gif
share/nagios/images/delete.gif
share/nagios/images/detail.gif
share/nagios/images/disabled.gif
share/nagios/images/down.gif
share/nagios/images/downtime.gif
share/nagios/images/empty.gif
share/nagios/images/enabled.gif
share/nagios/images/extinfo.gif
share/nagios/images/favicon.ico
share/nagios/images/flapping.gif
share/nagios/images/globe-support-150x150.png
share/nagios/images/graph.gif
share/nagios/images/greendot.gif
share/nagios/images/histogram.png
share/nagios/images/history.gif
share/nagios/images/hostevent.gif
share/nagios/images/info.png
share/nagios/images/left.gif
share/nagios/images/logofullsize.png
share/nagios/images/logos/aix.gd2
share/nagios/images/logos/aix.gif
share/nagios/images/logos/aix.jpg
share/nagios/images/logos/aix.png
share/nagios/images/logos/amiga.gd2
share/nagios/images/logos/amiga.gif
share/nagios/images/logos/amiga.jpg
share/nagios/images/logos/amiga.png
share/nagios/images/logos/apple.gd2
share/nagios/images/logos/apple.gif
share/nagios/images/logos/apple.jpg
share/nagios/images/logos/apple.png
share/nagios/images/logos/beos.gd2
share/nagios/images/logos/beos.gif
share/nagios/images/logos/beos.jpg
share/nagios/images/logos/beos.png
share/nagios/images/logos/bluetooth.png
share/nagios/images/logos/caldera.gd2
share/nagios/images/logos/caldera.gif
share/nagios/images/logos/caldera.jpg
share/nagios/images/logos/caldera.png
share/nagios/images/logos/cat1900.gd2
share/nagios/images/logos/cat2900.gd2
share/nagios/images/logos/cat5000.gd2
share/nagios/images/logos/database.gd2
share/nagios/images/logos/database.gif
share/nagios/images/logos/debian.gd2
share/nagios/images/logos/debian.gif
share/nagios/images/logos/debian.jpg
share/nagios/images/logos/debian.png
share/nagios/images/logos/desktop-server.gd2
share/nagios/images/logos/desktop-server.gif
share/nagios/images/logos/ethernet_card.png
share/nagios/images/logos/fax.gd2
share/nagios/images/logos/fax.gif
share/nagios/images/logos/firewall.gd2
share/nagios/images/logos/firewall.gif
share/nagios/images/logos/freebsd40.gd2
share/nagios/images/logos/freebsd40.gif
share/nagios/images/logos/freebsd40.jpg
share/nagios/images/logos/freebsd40.png
share/nagios/images/logos/globe.png
share/nagios/images/logos/graph.gif
share/nagios/images/logos/hp-printer40.gd2
share/nagios/images/logos/hp-printer40.gif
share/nagios/images/logos/hp-printer40.jpg
share/nagios/images/logos/hp-printer40.png
share/nagios/images/logos/hpux.gd2
share/nagios/images/logos/hpux.gif
share/nagios/images/logos/hpux.jpg
share/nagios/images/logos/hpux.png
share/nagios/images/logos/hub.gd2
share/nagios/images/logos/hub.gif
share/nagios/images/logos/internet.gd2
share/nagios/images/logos/internet.gif
share/nagios/images/logos/internet_device.png
share/nagios/images/logos/ip-pbx.gd2
share/nagios/images/logos/ip-pbx.gif
share/nagios/images/logos/irix.gd2
share/nagios/images/logos/irix.gif
share/nagios/images/logos/irix.jpg
share/nagios/images/logos/irix.png
share/nagios/images/logos/linux40.gd2
share/nagios/images/logos/linux40.gif
share/nagios/images/logos/linux40.jpg
share/nagios/images/logos/linux40.png
share/nagios/images/logos/logo.gd2
share/nagios/images/logos/mac40.gd2
share/nagios/images/logos/mac40.gif
share/nagios/images/logos/mac40.jpg
share/nagios/images/logos/mac40.png
share/nagios/images/logos/mainframe.gd2
share/nagios/images/logos/mainframe.gif
share/nagios/images/logos/mandrake.gd2
share/nagios/images/logos/mandrake.gif
share/nagios/images/logos/mandrake.jpg
share/nagios/images/logos/mandrake.png
share/nagios/images/logos/monitor.png
share/nagios/images/logos/nagios.gd2
share/nagios/images/logos/nagios.gif
share/nagios/images/logos/nagiosvrml.png
share/nagios/images/logos/netbsd40.gd2
share/nagios/images/logos/netbsd40.gif
share/nagios/images/logos/netbsd40.jpg
share/nagios/images/logos/netbsd40.png
share/nagios/images/logos/next.gd2
share/nagios/images/logos/next.gif
share/nagios/images/logos/next.jpg
share/nagios/images/logos/next.png
share/nagios/images/logos/ng-switch40.gd2
share/nagios/images/logos/ng-switch40.gif
share/nagios/images/logos/ng-switch40.jpg
share/nagios/images/logos/ng-switch40.png
share/nagios/images/logos/notebook.gd2
share/nagios/images/logos/notebook.gif
share/nagios/images/logos/novell40.gd2
share/nagios/images/logos/novell40.gif
share/nagios/images/logos/novell40.jpg
share/nagios/images/logos/novell40.png
share/nagios/images/logos/openbsd.gd2
share/nagios/images/logos/openbsd.gif
share/nagios/images/logos/openbsd.jpg
share/nagios/images/logos/openbsd.png
share/nagios/images/logos/printer.gd2
share/nagios/images/logos/printer.gif
share/nagios/images/logos/rack-server.gd2
share/nagios/images/logos/rack-server.gif
share/nagios/images/logos/redhat.gd2
share/nagios/images/logos/redhat.gif
share/nagios/images/logos/redhat.jpg
share/nagios/images/logos/redhat.png
share/nagios/images/logos/router.gd2
share/nagios/images/logos/router.gif
share/nagios/images/logos/router40.gd2
share/nagios/images/logos/router40.gif
share/nagios/images/logos/router40.jpg
share/nagios/images/logos/router40.png
share/nagios/images/logos/san.gd2
share/nagios/images/logos/san.gif
share/nagios/images/logos/satellite.png
share/nagios/images/logos/server.png
share/nagios/images/logos/signal.png
share/nagios/images/logos/slackware.gd2
share/nagios/images/logos/slackware.gif
share/nagios/images/logos/slackware.jpg
share/nagios/images/logos/slackware.png
share/nagios/images/logos/stampede.gd2
share/nagios/images/logos/stampede.gif
share/nagios/images/logos/stampede.jpg
share/nagios/images/logos/stampede.png
share/nagios/images/logos/station.gd2
share/nagios/images/logos/storm.gd2
share/nagios/images/logos/storm.gif
share/nagios/images/logos/storm.jpg
share/nagios/images/logos/storm.png
share/nagios/images/logos/sun40.gd2
share/nagios/images/logos/sun40.gif
share/nagios/images/logos/sun40.jpg
share/nagios/images/logos/sun40.png
share/nagios/images/logos/sunlogo.gd2
share/nagios/images/logos/sunlogo.gif
share/nagios/images/logos/sunlogo.jpg
share/nagios/images/logos/sunlogo.png
share/nagios/images/logos/switch.gd2
share/nagios/images/logos/switch.gif
share/nagios/images/logos/switch40.gd2
share/nagios/images/logos/switch40.gif
share/nagios/images/logos/switch40.jpg
share/nagios/images/logos/switch40.png
share/nagios/images/logos/thin-client.gd2
share/nagios/images/logos/thin-client.gif
share/nagios/images/logos/turbolinux.gd2
share/nagios/images/logos/turbolinux.gif
share/nagios/images/logos/turbolinux.jpg
share/nagios/images/logos/turbolinux.png
share/nagios/images/logos/ultrapenguin.gd2
share/nagios/images/logos/ultrapenguin.gif
share/nagios/images/logos/ultrapenguin.jpg
share/nagios/images/logos/ultrapenguin.png
share/nagios/images/logos/unicos.gd2
share/nagios/images/logos/unicos.gif
share/nagios/images/logos/unicos.jpg
share/nagios/images/logos/unicos.png
share/nagios/images/logos/unknown.gd2
share/nagios/images/logos/unknown.gif
share/nagios/images/logos/webcamera.png
share/nagios/images/logos/wifi.gd2
share/nagios/images/logos/wifi.gif
share/nagios/images/logos/wifi_modem.png
share/nagios/images/logos/win40.gd2
share/nagios/images/logos/win40.gif
share/nagios/images/logos/win40.jpg
share/nagios/images/logos/win40.png
share/nagios/images/logos/workstation.gd2
share/nagios/images/logos/workstation.gif
share/nagios/images/logos/workstation.png
share/nagios/images/logos/workstation_locked.png
share/nagios/images/logos/yellowdog.gd2
share/nagios/images/logos/yellowdog.gif
share/nagios/images/logos/yellowdog.jpg
share/nagios/images/logos/yellowdog.png
share/nagios/images/logrotate.png
share/nagios/images/ndisabled.gif
share/nagios/images/noack.gif
share/nagios/images/notes.gif
share/nagios/images/notify.gif
share/nagios/images/orangedot.gif
share/nagios/images/passiveonly.gif
share/nagios/images/recovery.png
share/nagios/images/redudancy.png
share/nagios/images/redundancy.png
share/nagios/images/restart.gif
share/nagios/images/right.gif
share/nagios/images/sblogo.png
share/nagios/images/serviceevent.gif
share/nagios/images/sflogo.png
share/nagios/images/splunk1.gif
share/nagios/images/splunk2.gif
share/nagios/images/start.gif
share/nagios/images/status.gif
share/nagios/images/status2.gif
share/nagios/images/status3.gif
share/nagios/images/status4.gif
share/nagios/images/stop.gif
share/nagios/images/tacdisabled.jpg
share/nagios/images/tacdisabled.png
share/nagios/images/tacenabled.jpg
share/nagios/images/tacenabled.png
share/nagios/images/thermcrit.png
share/nagios/images/thermok.png
share/nagios/images/thermwarn.png
share/nagios/images/trends.gif
share/nagios/images/trendshost.png
share/nagios/images/trendssvc.png
share/nagios/images/unknown.png
share/nagios/images/up.gif
share/nagios/images/warning.png
share/nagios/images/weblogo1.png
share/nagios/images/zoom1.gif
share/nagios/images/zoom2.gif
share/nagios/includes/jquery-1.7.1.min.js
share/nagios/includes/rss/extlib/Snoopy.class.inc
share/nagios/includes/rss/rss_cache.inc
share/nagios/includes/rss/rss_fetch.inc
share/nagios/includes/rss/rss_parse.inc
share/nagios/includes/rss/rss_utils.inc
share/nagios/includes/utils.inc.php
share/nagios/index.php
share/nagios/js/jquery-1.7.1.min.js
share/nagios/main.php
share/nagios/robots.txt
share/nagios/rss-corefeed.php
share/nagios/rss-newsfeed.php
share/nagios/side.php
share/nagios/stylesheets/avail.css
share/nagios/stylesheets/checksanity.css
share/nagios/stylesheets/cmd.css
share/nagios/stylesheets/common.css
share/nagios/stylesheets/config.css
share/nagios/stylesheets/extinfo.css
share/nagios/stylesheets/histogram.css
share/nagios/stylesheets/history.css
share/nagios/stylesheets/ministatus.css
share/nagios/stylesheets/notifications.css
share/nagios/stylesheets/outages.css
share/nagios/stylesheets/showlog.css
share/nagios/stylesheets/status.css
share/nagios/stylesheets/statusmap.css
share/nagios/stylesheets/summary.css
share/nagios/stylesheets/tac.css
share/nagios/stylesheets/trends.css
@pkgdir share/nagios/ssi
@pkgdir etc/nagios/objects
