# net/nagios-base/PLIST
@comment $NetBSD: PLIST,v 1.7 2009/06/14 21:00:04 joerg Exp $
bin/nagios
bin/nagiostats
bin/p1.pl
libexec/nagios/cgi-bin/avail.cgi
libexec/nagios/cgi-bin/cmd.cgi
libexec/nagios/cgi-bin/config.cgi
libexec/nagios/cgi-bin/extinfo.cgi
libexec/nagios/cgi-bin/histogram.cgi
libexec/nagios/cgi-bin/history.cgi
libexec/nagios/cgi-bin/notifications.cgi
libexec/nagios/cgi-bin/outages.cgi
libexec/nagios/cgi-bin/showlog.cgi
libexec/nagios/cgi-bin/status.cgi
libexec/nagios/cgi-bin/statusmap.cgi
libexec/nagios/cgi-bin/statuswml.cgi
libexec/nagios/cgi-bin/statuswrl.cgi
libexec/nagios/cgi-bin/summary.cgi
libexec/nagios/cgi-bin/tac.cgi
libexec/nagios/cgi-bin/trends.cgi
share/examples/nagios/cgi.cfg-sample
share/examples/nagios/nagios.cfg-sample
share/examples/nagios/nagios.conf
share/examples/nagios/objects/commands.cfg-sample
share/examples/nagios/objects/contacts.cfg-sample
share/examples/nagios/objects/localhost.cfg-sample
share/examples/nagios/objects/printer.cfg-sample
share/examples/nagios/objects/switch.cfg-sample
share/examples/nagios/objects/templates.cfg-sample
share/examples/nagios/objects/timeperiods.cfg-sample
share/examples/nagios/objects/windows.cfg-sample
share/examples/nagios/resource.cfg-sample
share/examples/rc.d/nagios
share/nagios/config.inc.php
share/nagios/contexthelp/A1.html
share/nagios/contexthelp/A2.html
share/nagios/contexthelp/A3.html
share/nagios/contexthelp/A4.html
share/nagios/contexthelp/A5.html
share/nagios/contexthelp/A6.html
share/nagios/contexthelp/A7.html
share/nagios/contexthelp/B1.html
share/nagios/contexthelp/C1.html
share/nagios/contexthelp/D1.html
share/nagios/contexthelp/E1.html
share/nagios/contexthelp/F1.html
share/nagios/contexthelp/G1.html
share/nagios/contexthelp/G2.html
share/nagios/contexthelp/G3.html
share/nagios/contexthelp/G4.html
share/nagios/contexthelp/G5.html
share/nagios/contexthelp/G6.html
share/nagios/contexthelp/H1.html
share/nagios/contexthelp/H2.html
share/nagios/contexthelp/H3.html
share/nagios/contexthelp/H4.html
share/nagios/contexthelp/H5.html
share/nagios/contexthelp/H6.html
share/nagios/contexthelp/H7.html
share/nagios/contexthelp/H8.html
share/nagios/contexthelp/I1.html
share/nagios/contexthelp/I2.html
share/nagios/contexthelp/I3.html
share/nagios/contexthelp/I4.html
share/nagios/contexthelp/I5.html
share/nagios/contexthelp/I6.html
share/nagios/contexthelp/I7.html
share/nagios/contexthelp/I8.html
share/nagios/contexthelp/I9.html
share/nagios/contexthelp/J1.html
share/nagios/contexthelp/K1.html
share/nagios/contexthelp/L1.html
share/nagios/contexthelp/L10.html
share/nagios/contexthelp/L11.html
share/nagios/contexthelp/L12.html
share/nagios/contexthelp/L13.html
share/nagios/contexthelp/L2.html
share/nagios/contexthelp/L3.html
share/nagios/contexthelp/L4.html
share/nagios/contexthelp/L5.html
share/nagios/contexthelp/L6.html
share/nagios/contexthelp/L7.html
share/nagios/contexthelp/L8.html
share/nagios/contexthelp/L9.html
share/nagios/contexthelp/M1.html
share/nagios/contexthelp/M2.html
share/nagios/contexthelp/M3.html
share/nagios/contexthelp/M4.html
share/nagios/contexthelp/M5.html
share/nagios/contexthelp/M6.html
share/nagios/contexthelp/N1.html
share/nagios/contexthelp/N2.html
share/nagios/contexthelp/N3.html
share/nagios/contexthelp/N4.html
share/nagios/contexthelp/N5.html
share/nagios/contexthelp/N6.html
share/nagios/contexthelp/N7.html
share/nagios/docs/about.html
share/nagios/docs/activechecks.html
share/nagios/docs/adaptive.html
share/nagios/docs/addons.html
share/nagios/docs/beginners.html
share/nagios/docs/cachedchecks.html
share/nagios/docs/cgiauth.html
share/nagios/docs/cgiincludes.html
share/nagios/docs/cgis.html
share/nagios/docs/cgisecurity.html
share/nagios/docs/checkscheduling.html
share/nagios/docs/clusters.html
share/nagios/docs/config.html
share/nagios/docs/configcgi.html
share/nagios/docs/configmain.html
share/nagios/docs/configobject.html
share/nagios/docs/customobjectvars.html
share/nagios/docs/dependencies.html
share/nagios/docs/dependencychecks.html
share/nagios/docs/distributed.html
share/nagios/docs/downtime.html
share/nagios/docs/embeddedperl.html
share/nagios/docs/epnplugins.html
share/nagios/docs/escalations.html
share/nagios/docs/eventhandlers.html
share/nagios/docs/extcommands.html
share/nagios/docs/faststartup.html
share/nagios/docs/flapping.html
share/nagios/docs/freshness.html
share/nagios/docs/funstuff.html
share/nagios/docs/hostchecks.html
share/nagios/docs/images/activechecks.png
share/nagios/docs/images/cachedcheckgraphs.png
share/nagios/docs/images/cachedchecks.png
share/nagios/docs/images/cachedchecks1.png
share/nagios/docs/images/cgi-avail-a.png
share/nagios/docs/images/cgi-avail-b.png
share/nagios/docs/images/cgi-cmd.png
share/nagios/docs/images/cgi-config.png
share/nagios/docs/images/cgi-extinfo-a.png
share/nagios/docs/images/cgi-extinfo-b.png
share/nagios/docs/images/cgi-extinfo-c.png
share/nagios/docs/images/cgi-extinfo-d.png
share/nagios/docs/images/cgi-histogram.png
share/nagios/docs/images/cgi-history.png
share/nagios/docs/images/cgi-notifications.png
share/nagios/docs/images/cgi-outages.png
share/nagios/docs/images/cgi-showlog.png
share/nagios/docs/images/cgi-status-a.png
share/nagios/docs/images/cgi-status-b.png
share/nagios/docs/images/cgi-status-c.png
share/nagios/docs/images/cgi-status-d.png
share/nagios/docs/images/cgi-statusmap.png
share/nagios/docs/images/cgi-statuswml.png
share/nagios/docs/images/cgi-statuswrl.png
share/nagios/docs/images/cgi-summary.png
share/nagios/docs/images/cgi-tac.png
share/nagios/docs/images/cgi-trends.png
share/nagios/docs/images/checkmark.png
share/nagios/docs/images/checktiming.png
share/nagios/docs/images/configoverview.png
share/nagios/docs/images/distributed.png
share/nagios/docs/images/downtime.png
share/nagios/docs/images/epn.png
share/nagios/docs/images/eventhandlers.png
share/nagios/docs/images/externalcommands.png
share/nagios/docs/images/fast-startup.png
share/nagios/docs/images/fast-startup1.png
share/nagios/docs/images/fast-startup2.png
share/nagios/docs/images/freshness.png
share/nagios/docs/images/host-dependencies.png
share/nagios/docs/images/important.gif
share/nagios/docs/images/integrationoverview.png
share/nagios/docs/images/interleaved1.png
share/nagios/docs/images/interleaved2.png
share/nagios/docs/images/interleaved3.png
share/nagios/docs/images/logofullsize.png
share/nagios/docs/images/monitoring-printers.png
share/nagios/docs/images/monitoring-routers.png
share/nagios/docs/images/monitoring-windows.png
share/nagios/docs/images/mrtg-activehostchecks.png
share/nagios/docs/images/mrtg-activelychecked.png
share/nagios/docs/images/mrtg-activeservicechecks.png
share/nagios/docs/images/mrtg-cachedchecks.png
share/nagios/docs/images/mrtg-commandbuffers.png
share/nagios/docs/images/mrtg-externalcommands.png
share/nagios/docs/images/mrtg-hostperfstats.png
share/nagios/docs/images/mrtg-hoststatechange.png
share/nagios/docs/images/mrtg-passivechecks.png
share/nagios/docs/images/mrtg-passivelychecked.png
share/nagios/docs/images/mrtg-serviceperfstats.png
share/nagios/docs/images/mrtg-servicestatechange.png
share/nagios/docs/images/multiple-templates1.png
share/nagios/docs/images/multiple-templates2.png
share/nagios/docs/images/nagios.jpg
share/nagios/docs/images/ndoutils.png
share/nagios/docs/images/network-outage1.png
share/nagios/docs/images/network-outage2.png
share/nagios/docs/images/noninterleaved1.png
share/nagios/docs/images/noninterleaved2.png
share/nagios/docs/images/note.gif
share/nagios/docs/images/nrpe.png
share/nagios/docs/images/nsca.png
share/nagios/docs/images/nscpp.png
share/nagios/docs/images/objects-commands.png
share/nagios/docs/images/objects-contacts.png
share/nagios/docs/images/objects-hosts.png
share/nagios/docs/images/objects-services.png
share/nagios/docs/images/objects-timeperiods.png
share/nagios/docs/images/passivechecks.png
share/nagios/docs/images/passivehosttranslation.png
share/nagios/docs/images/plugins.png
share/nagios/docs/images/predictive-dependency-checks.png
share/nagios/docs/images/printer.png
share/nagios/docs/images/reachability1.png
share/nagios/docs/images/reachability2.png
share/nagios/docs/images/reachability3.png
share/nagios/docs/images/reachability4.png
share/nagios/docs/images/redudancy.png
share/nagios/docs/images/redundancy.png
share/nagios/docs/images/security.png
share/nagios/docs/images/security1.png
share/nagios/docs/images/security2.png
share/nagios/docs/images/security3.png
share/nagios/docs/images/seealso.gif
share/nagios/docs/images/service-dependencies.png
share/nagios/docs/images/statetransitions.png
share/nagios/docs/images/statetransitions2.png
share/nagios/docs/images/stoprestart.png
share/nagios/docs/images/switch.png
share/nagios/docs/images/tcpwrappers.png
share/nagios/docs/images/tip.gif
share/nagios/docs/images/tuning.png
share/nagios/docs/images/upto.gif
share/nagios/docs/index.html
share/nagios/docs/int-snmptrap.html
share/nagios/docs/int-tcpwrappers.html
share/nagios/docs/integration.html
share/nagios/docs/largeinstalltweaks.html
share/nagios/docs/macrolist.html
share/nagios/docs/macros.html
share/nagios/docs/monitoring-linux.html
share/nagios/docs/monitoring-netware.html
share/nagios/docs/monitoring-printers.html
share/nagios/docs/monitoring-publicservices.html
share/nagios/docs/monitoring-routers.html
share/nagios/docs/monitoring-windows.html
share/nagios/docs/mrtggraphs.html
share/nagios/docs/nagiostats.html
share/nagios/docs/networkreachability.html
share/nagios/docs/notifications.html
share/nagios/docs/objectdefinitions.html
share/nagios/docs/objectinheritance.html
share/nagios/docs/objecttricks.html
share/nagios/docs/oncallrotation.html
share/nagios/docs/passivechecks.html
share/nagios/docs/passivestatetranslation.html
share/nagios/docs/perfdata.html
share/nagios/docs/pluginapi.html
share/nagios/docs/plugins.html
share/nagios/docs/quickstart-fedora.html
share/nagios/docs/quickstart-opensuse.html
share/nagios/docs/quickstart-ubuntu.html
share/nagios/docs/quickstart.html
share/nagios/docs/redundancy.html
share/nagios/docs/robots.txt
share/nagios/docs/security.html
share/nagios/docs/servicechecks.html
share/nagios/docs/stalking.html
share/nagios/docs/startstop.html
share/nagios/docs/statetypes.html
share/nagios/docs/timeperiods.html
share/nagios/docs/toc.html
share/nagios/docs/tuning.html
share/nagios/docs/upgrading.html
share/nagios/docs/verifyconfig.html
share/nagios/docs/volatileservices.html
share/nagios/docs/whatsnew.html
share/nagios/images/NagiosEnterprises-whitebg-112x46.png
share/nagios/images/ack.gif
share/nagios/images/action.gif
share/nagios/images/command.png
share/nagios/images/comment.gif
share/nagios/images/contexthelp1.gif
share/nagios/images/contexthelp2.gif
share/nagios/images/critical.png
share/nagios/images/delay.gif
share/nagios/images/delete.gif
share/nagios/images/detail.gif
share/nagios/images/disabled.gif
share/nagios/images/down.gif
share/nagios/images/downtime.gif
share/nagios/images/empty.gif
share/nagios/images/enabled.gif
share/nagios/images/extinfo.gif
share/nagios/images/favicon.ico
share/nagios/images/flapping.gif
share/nagios/images/globe-support-150x150.png
share/nagios/images/greendot.gif
share/nagios/images/histogram.png
share/nagios/images/history.gif
share/nagios/images/hostevent.gif
share/nagios/images/info.png
share/nagios/images/left.gif
share/nagios/images/logofullsize.png
share/nagios/images/logos/nagios.gd2
share/nagios/images/logos/nagios.gif
share/nagios/images/logos/nagiosvrml.png
share/nagios/images/logos/unknown.gd2
share/nagios/images/logos/unknown.gif
share/nagios/images/logrotate.png
share/nagios/images/ndisabled.gif
share/nagios/images/noack.gif
share/nagios/images/notes.gif
share/nagios/images/notify.gif
share/nagios/images/orangedot.gif
share/nagios/images/passiveonly.gif
share/nagios/images/recovery.png
share/nagios/images/redudancy.png
share/nagios/images/redundancy.png
share/nagios/images/restart.gif
share/nagios/images/right.gif
share/nagios/images/sblogo.png
share/nagios/images/serviceevent.gif
share/nagios/images/splunk1.gif
share/nagios/images/splunk2.gif
share/nagios/images/start.gif
share/nagios/images/status.gif
share/nagios/images/status2.gif
share/nagios/images/status3.gif
share/nagios/images/status4.gif
share/nagios/images/stop.gif
share/nagios/images/tacdisabled.jpg
share/nagios/images/tacdisabled.png
share/nagios/images/tacenabled.jpg
share/nagios/images/tacenabled.png
share/nagios/images/thermcrit.png
share/nagios/images/thermok.png
share/nagios/images/thermwarn.png
share/nagios/images/trends.gif
share/nagios/images/trendshost.png
share/nagios/images/trendssvc.png
share/nagios/images/unknown.png
share/nagios/images/up.gif
share/nagios/images/warning.png
share/nagios/images/weblogo1.png
share/nagios/images/zoom1.gif
share/nagios/images/zoom2.gif
share/nagios/includes/utils.inc.php
share/nagios/index.php
share/nagios/main.php
share/nagios/robots.txt
share/nagios/side.php
share/nagios/stylesheets/avail.css
share/nagios/stylesheets/checksanity.css
share/nagios/stylesheets/cmd.css
share/nagios/stylesheets/common.css
share/nagios/stylesheets/config.css
share/nagios/stylesheets/extinfo.css
share/nagios/stylesheets/histogram.css
share/nagios/stylesheets/history.css
share/nagios/stylesheets/ministatus.css
share/nagios/stylesheets/notifications.css
share/nagios/stylesheets/outages.css
share/nagios/stylesheets/showlog.css
share/nagios/stylesheets/status.css
share/nagios/stylesheets/statusmap.css
share/nagios/stylesheets/summary.css
share/nagios/stylesheets/tac.css
share/nagios/stylesheets/trends.css
@pkgdir share/nagios/ssi
