# net/hobbitclient/PLIST
@comment $NetBSD: PLIST,v 1.3 2009/06/14 18:09:30 joerg Exp $
@comment $NetBSD: PLIST,v 1.3 2009/06/14 18:09:30 joerg Exp $
hobbit/client/bin/bb
hobbit/client/bin/bbcmd
hobbit/client/bin/bbdigest
hobbit/client/bin/bbhostgrep
hobbit/client/bin/bbhostshow
hobbit/client/bin/clientupdate
hobbit/client/bin/hobbitlaunch
hobbit/client/bin/logfetch
hobbit/client/bin/msgcache
hobbit/client/bin/${LOWER_OPSYS}-meminfo
hobbit/client/bin/orcahobbit
hobbit/client/etc/clientlaunch.cfg
hobbit/client/etc/hobbitclient.cfg
hobbit/client/etc/localclient.cfg
hobbit/client/bin/hobbitclient-aix.sh
hobbit/client/bin/hobbitclient-darwin.sh
hobbit/client/bin/hobbitclient-freebsd.sh
hobbit/client/bin/hobbitclient-hp-ux.sh
hobbit/client/bin/hobbitclient-irix.sh
hobbit/client/bin/hobbitclient-linux.sh
hobbit/client/bin/hobbitclient-netbsd.sh
hobbit/client/bin/hobbitclient-openbsd.sh
hobbit/client/bin/hobbitclient-osf1.sh
hobbit/client/bin/hobbitclient-sco_sv.sh
hobbit/client/bin/hobbitclient-sunos.sh
hobbit/client/bin/hobbitclient.sh
hobbit/client/runclient.sh
share/examples/rc.d/hobbitclient
@exec ${MKDIR} %D/hobbit/client/tmp
@exec ${MKDIR} %D/hobbit/client/logs
@exec ${MKDIR} %D/hobbit/client/ext
