# net/socks5/PLIST
@comment $NetBSD: PLIST,v 1.6 2009/06/12 19:09:38 rillig Exp $
bin/rarchie
bin/rfinger
bin/rftp
bin/rping
bin/rtelnet
bin/rtraceroute
bin/runsocks
bin/rwhois
bin/socks5
bin/stopsocks
include/socks.h
lib/libsocks5.a
lib/libsocks5_sh.so
lib/libsocks5_sh.so.0
lib/libsocks5_sh.so.0.0
man/man1/runsocks.1
man/man1/socks5.1
man/man1/socks5_clients.1
man/man1/stopsocks.1
man/man5/libsocks5.conf.5
man/man5/socks5.conf.5
man/man5/socks5.passwd.5
share/examples/rc.d/socks5
share/examples/socks5/README
share/examples/socks5/multiple_server.example
share/examples/socks5/server_2_multiple_server.example
share/examples/socks5/single_homed_host.example
share/examples/socks5/socks5.conf.dualhomed
share/examples/socks5/socks5.conf.gssapi
share/examples/socks5/socks5.conf.multipleservers
share/examples/socks5/socks5.conf.server2server
share/examples/socks5/socks5.conf.singlehomed
@dirrm share/examples/socks5
