# net/dhid/PLIST
@comment $NetBSD: PLIST,v 1.4 2009/06/14 18:09:25 joerg Exp $
sbin/dhid
sbin/dhid-genkeys
share/doc/dhid/COPYRIGHT
share/doc/dhid/README
share/doc/dhid/AUTHORS
share/doc/dhid/WHATSNEW
share/doc/dhid/INSTALL
share/examples/rc.d/dhid
