# net/dhid/PLIST
@comment $NetBSD: PLIST,v 1.1 2001/10/31 22:55:12 zuntum Exp $
sbin/dhid
sbin/dhid-genkeys
etc/rc.d/dhid
share/doc/dhid/COPYRIGHT
share/doc/dhid/README
share/doc/dhid/AUTHORS
share/doc/dhid/WHATSNEW
share/doc/dhid/INSTALL
@dirrm share/doc/dhid
