# sysutils/ipa/PLIST
@comment $NetBSD: PLIST,v 1.2 2007/01/23 16:59:51 reed Exp $
bin/ipa
bin/ipactl
bin/ipastat
include/ipa_mod.h
man/man3/ipa_mod.3
man/man5/ipa.conf.5
man/man5/ipastat.conf.5
man/man8/ipa.8
man/man8/ipactl.8
man/man8/ipastat.8
man/ru/man3/ipa_mod.3
man/ru/man5/ipa.conf.5
man/ru/man5/ipastat.conf.5
man/ru/man8/ipa.8
man/ru/man8/ipactl.8
man/ru/man8/ipastat.8
share/examples/rc.d/ipa
