# net/proftpd/PLIST
@comment $NetBSD: PLIST,v 1.10.8.1 2006/12/18 14:28:32 ghen Exp $
bin/ftpcount
bin/ftpdctl
bin/ftptop
bin/ftpwho
include/proftpd/auth.h
include/proftpd/bindings.h
include/proftpd/buildstamp.h
include/proftpd/child.h
include/proftpd/class.h
include/proftpd/compat.h
include/proftpd/conf.h
include/proftpd/config.h
include/proftpd/ctrls.h
include/proftpd/data.h
include/proftpd/default_paths.h
include/proftpd/dirtree.h
include/proftpd/display.h
include/proftpd/event.h
include/proftpd/feat.h
include/proftpd/fsio.h
include/proftpd/ftp.h
include/proftpd/glibc-glob.h
include/proftpd/help.h
include/proftpd/ident.h
include/proftpd/inet.h
include/proftpd/lastlog.h
include/proftpd/libsupp.h
include/proftpd/log.h
include/proftpd/mkhome.h
include/proftpd/mod_ctrls.h
include/proftpd/modules.h
include/proftpd/netacl.h
include/proftpd/netaddr.h
include/proftpd/netio.h
include/proftpd/options.h
include/proftpd/parser.h
include/proftpd/pool.h
include/proftpd/pr-syslog.h
include/proftpd/privs.h
include/proftpd/proftpd.h
include/proftpd/regexp.h
include/proftpd/response.h
include/proftpd/scoreboard.h
include/proftpd/sets.h
include/proftpd/support.h
include/proftpd/table.h
include/proftpd/timers.h
include/proftpd/trace.h
include/proftpd/utf8.h
include/proftpd/var.h
include/proftpd/version.h
include/proftpd/xferlog.h
man/man1/ftpcount.1
man/man1/ftptop.1
man/man1/ftpwho.1
man/man5/xferlog.5
man/man8/ftpdctl.8
man/man8/ftpshut.8
man/man8/proftpd.8
sbin/ftpshut
sbin/in.proftpd
sbin/proftpd
share/doc/proftpd/Configuration.html
share/doc/proftpd/NEWS
share/doc/proftpd/README
share/doc/proftpd/README.IPv6
share/doc/proftpd/README.LDAP
share/doc/proftpd/README.PAM
share/doc/proftpd/README.capabilities
share/doc/proftpd/README.classes
share/doc/proftpd/README.controls
share/doc/proftpd/README.facl
share/doc/proftpd/README.modules
share/doc/proftpd/faq.html
share/doc/proftpd/license.txt
share/examples/proftpd/anonymous.conf
share/examples/proftpd/basic.conf
share/examples/proftpd/complex-virtual.conf
share/examples/proftpd/mod_sql.conf
share/examples/proftpd/virtual.conf
share/examples/rc.d/proftpd
@dirrm share/doc/proftpd
@dirrm share/examples/proftpd
@dirrm include/proftpd
