# net/proftpd/PLIST
@comment $NetBSD: PLIST,v 1.8 2005/05/02 20:34:03 reed Exp $
bin/ftpcount
bin/ftpwho
bin/ftpdctl
bin/ftptop
man/man1/ftpcount.1
man/man1/ftptop.1
man/man1/ftpwho.1
man/man5/xferlog.5
man/man8/ftpdctl.8
man/man8/ftpshut.8
man/man8/proftpd.8
sbin/ftpshut
sbin/in.proftpd
sbin/proftpd
share/doc/proftpd/Configuration.html
share/doc/proftpd/NEWS
share/doc/proftpd/README
share/doc/proftpd/README.IPv6
share/doc/proftpd/README.LDAP
share/doc/proftpd/README.PAM
share/doc/proftpd/README.modules
share/doc/proftpd/faq.html
share/doc/proftpd/license.txt
share/examples/proftpd/anonymous.conf
share/examples/proftpd/basic.conf
share/examples/proftpd/complex-virtual.conf
share/examples/proftpd/mod_sql.conf
share/examples/proftpd/virtual.conf
share/examples/rc.d/proftpd
@dirrm share/doc/proftpd
@dirrm share/examples/proftpd
