# net/proftpd/PLIST
@comment $NetBSD: PLIST,v 1.23.16.1 2013/03/17 11:45:23 tron Exp $
bin/ftpasswd
bin/ftpcount
bin/ftpdctl
bin/ftpmail
bin/ftpquota
bin/ftptop
bin/ftpwho
bin/prxs
include/proftpd/auth.h
include/proftpd/bindings.h
include/proftpd/buildstamp.h
include/proftpd/child.h
include/proftpd/class.h
include/proftpd/cmd.h
include/proftpd/compat.h
include/proftpd/conf.h
include/proftpd/config.h
include/proftpd/ctrls.h
include/proftpd/data.h
include/proftpd/default_paths.h
include/proftpd/dirtree.h
include/proftpd/display.h
include/proftpd/encode.h
include/proftpd/env.h
include/proftpd/event.h
include/proftpd/expr.h
include/proftpd/feat.h
include/proftpd/filter.h
include/proftpd/fsio.h
include/proftpd/ftp.h
include/proftpd/glibc-glob.h
include/proftpd/help.h
include/proftpd/inet.h
include/proftpd/libsupp.h
include/proftpd/log.h
include/proftpd/memcache.h
include/proftpd/mkhome.h
include/proftpd/mod_ctrls.h
${PLIST.quota}include/proftpd/mod_quotatab.h
${PLIST.sql}include/proftpd/mod_sql.h
${PLIST.tls}include/proftpd/mod_tls.h
include/proftpd/modules.h
include/proftpd/netacl.h
include/proftpd/netaddr.h
include/proftpd/netio.h
include/proftpd/options.h
include/proftpd/parser.h
include/proftpd/pidfile.h
include/proftpd/pool.h
include/proftpd/pr-syslog.h
include/proftpd/privs.h
include/proftpd/proctitle.h
include/proftpd/proftpd.h
include/proftpd/regexp.h
include/proftpd/response.h
include/proftpd/scoreboard.h
include/proftpd/session.h
include/proftpd/sets.h
include/proftpd/stash.h
include/proftpd/str.h
include/proftpd/support.h
include/proftpd/table.h
include/proftpd/throttle.h
include/proftpd/timers.h
include/proftpd/tpl.h
include/proftpd/trace.h
include/proftpd/var.h
include/proftpd/version.h
include/proftpd/xferlog.h
lib/pkgconfig/proftpd.pc
man/man1/ftpasswd.1
man/man1/ftpcount.1
man/man1/ftpmail.1
man/man1/ftpquota.1
man/man1/ftptop.1
man/man1/ftpwho.1
man/man5/xferlog.5
man/man8/ftpdctl.8
man/man8/ftpscrub.8
man/man8/ftpshut.8
man/man8/proftpd.8
sbin/ftpscrub
sbin/ftpshut
sbin/in.proftpd
sbin/proftpd
share/doc/proftpd/Configuration.html
share/doc/proftpd/NEWS
share/doc/proftpd/README
share/doc/proftpd/README.IPv6
share/doc/proftpd/README.LDAP
share/doc/proftpd/README.PAM
share/doc/proftpd/README.capabilities
share/doc/proftpd/README.classes
share/doc/proftpd/README.controls
share/doc/proftpd/README.facl
share/doc/proftpd/README.modules
share/doc/proftpd/faq.html
share/doc/proftpd/license.txt
share/examples/proftpd/anonymous.conf
share/examples/proftpd/basic.conf
share/examples/proftpd/complex-virtual.conf
share/examples/proftpd/ftpasswd
share/examples/proftpd/mod_sql.conf
share/examples/proftpd/virtual.conf
share/examples/rc.d/proftpd
share/locale/bg_BG/LC_MESSAGES/proftpd.mo
share/locale/en_US/LC_MESSAGES/proftpd.mo
share/locale/es_ES/LC_MESSAGES/proftpd.mo
share/locale/fr_FR/LC_MESSAGES/proftpd.mo
share/locale/it_IT/LC_MESSAGES/proftpd.mo
share/locale/ja_JP/LC_MESSAGES/proftpd.mo
share/locale/ko_KR/LC_MESSAGES/proftpd.mo
share/locale/ru_RU/LC_MESSAGES/proftpd.mo
share/locale/zh_CN/LC_MESSAGES/proftpd.mo
share/locale/zh_TW/LC_MESSAGES/proftpd.mo
@pkgdir lib/proftpd
