# security/openct/PLIST
@comment $NetBSD: PLIST,v 1.3 2009/06/14 18:13:34 joerg Exp $
bin/openct-tool
include/openct/apdu.h
include/openct/buffer.h
include/openct/conf.h
include/openct/ctapi.h
include/openct/device.h
include/openct/driver.h
include/openct/error.h
include/openct/ifd.h
include/openct/logging.h
include/openct/openct.h
include/openct/path.h
include/openct/protocol.h
include/openct/server.h
include/openct/socket.h
include/openct/tlv.h
include/openct/types.h
lib/libopenct.la
lib/libopenctapi.la
${PLIST.pcsc-lite}lib/openct-ifd.la
lib/pkgconfig/libopenct.pc
man/man1/openct-tool.1
sbin/ifdhandler
sbin/ifdproxy
sbin/openct-control
share/doc/openct/NEWS
share/doc/openct/README
share/doc/openct/api/annotated.html
share/doc/openct/api/apdu_8h-source.html
share/doc/openct/api/buffer_8h-source.html
share/doc/openct/api/classes.html
share/doc/openct/api/conf_8h-source.html
share/doc/openct/api/device_8h-source.html
share/doc/openct/api/dir_c35229a911af2f480aef06e8035e7826.html
share/doc/openct/api/dirs.html
share/doc/openct/api/doxygen.css
share/doc/openct/api/doxygen.png
share/doc/openct/api/driver_8h-source.html
share/doc/openct/api/driver_8h.html
share/doc/openct/api/error_8h-source.html
share/doc/openct/api/files.html
share/doc/openct/api/functions.html
share/doc/openct/api/functions_vars.html
share/doc/openct/api/ifd_8h-source.html
share/doc/openct/api/index.html
share/doc/openct/api/logging_8h-source.html
share/doc/openct/api/openct_8h-source.html
share/doc/openct/api/opensc-logo.gif
share/doc/openct/api/path_8h-source.html
share/doc/openct/api/protocol_8h-source.html
share/doc/openct/api/server_8h-source.html
share/doc/openct/api/socket_8h-source.html
share/doc/openct/api/structifd__driver__ops.html
share/doc/openct/api/tab_b.gif
share/doc/openct/api/tab_l.gif
share/doc/openct/api/tab_r.gif
share/doc/openct/api/tabs.css
share/doc/openct/api/tlv_8h-source.html
share/doc/openct/api/types_8h-source.html
share/examples/openct/openct.conf
share/examples/openct/reader.conf
share/examples/rc.d/openct
