# mail/imapproxy/PLIST
@comment $NetBSD: PLIST,v 1.3 2010/03/22 08:10:36 manu Exp $
sbin/in.imapproxyd
sbin/pimpstat
share/doc/imapproxy/COPYING
share/doc/imapproxy/ChangeLog
share/doc/imapproxy/README
share/doc/imapproxy/README.known_issues
share/doc/imapproxy/README.ssl
share/examples/imapproxy/imapproxy.conf
share/examples/rc.d/imapproxy
