# www/lighttpd/PLIST
@comment $NetBSD: PLIST,v 1.9 2010/02/08 14:47:54 joerg Exp $
lib/lighttpd/mod_access.la
lib/lighttpd/mod_accesslog.la
lib/lighttpd/mod_alias.la
lib/lighttpd/mod_auth.la
lib/lighttpd/mod_cgi.la
lib/lighttpd/mod_cml.la
lib/lighttpd/mod_compress.la
lib/lighttpd/mod_dirlisting.la
lib/lighttpd/mod_evasive.la
lib/lighttpd/mod_evhost.la
lib/lighttpd/mod_expire.la
lib/lighttpd/mod_extforward.la
lib/lighttpd/mod_fastcgi.la
lib/lighttpd/mod_flv_streaming.la
lib/lighttpd/mod_indexfile.la
lib/lighttpd/mod_magnet.la
lib/lighttpd/mod_mysql_vhost.la
lib/lighttpd/mod_proxy.la
lib/lighttpd/mod_redirect.la
lib/lighttpd/mod_rewrite.la
lib/lighttpd/mod_rrdtool.la
lib/lighttpd/mod_scgi.la
lib/lighttpd/mod_secdownload.la
lib/lighttpd/mod_setenv.la
lib/lighttpd/mod_simple_vhost.la
lib/lighttpd/mod_ssi.la
lib/lighttpd/mod_staticfile.la
lib/lighttpd/mod_status.la
lib/lighttpd/mod_trigger_b4_dl.la
lib/lighttpd/mod_userdir.la
lib/lighttpd/mod_usertrack.la
lib/lighttpd/mod_webdav.la
man/man8/lighttpd.8
sbin/lighttpd
sbin/lighttpd-angel
share/doc/lighttpd/AUTHORS
share/doc/lighttpd/COPYING
share/doc/lighttpd/NEWS
share/doc/lighttpd/README
share/doc/lighttpd/access.txt
share/doc/lighttpd/accesslog.txt
share/doc/lighttpd/alias.txt
share/doc/lighttpd/authentication.txt
share/doc/lighttpd/cgi.txt
share/doc/lighttpd/cml.txt
share/doc/lighttpd/compress.txt
share/doc/lighttpd/configuration.txt
share/doc/lighttpd/dirlisting.txt
share/doc/lighttpd/evhost.txt
share/doc/lighttpd/expire.txt
share/doc/lighttpd/extforward.txt
share/doc/lighttpd/fastcgi-state.dot
share/doc/lighttpd/fastcgi-state.txt
share/doc/lighttpd/fastcgi.txt
share/doc/lighttpd/features.txt
share/doc/lighttpd/magnet.txt
share/doc/lighttpd/mysqlvhost.txt
share/doc/lighttpd/newstyle.css
share/doc/lighttpd/oldstyle.css
share/doc/lighttpd/performance.txt
share/doc/lighttpd/plugins.txt
share/doc/lighttpd/proxy.txt
share/doc/lighttpd/redirect.txt
share/doc/lighttpd/rewrite.txt
share/doc/lighttpd/rrdtool.txt
share/doc/lighttpd/scgi.txt
share/doc/lighttpd/secdownload.txt
share/doc/lighttpd/security.txt
share/doc/lighttpd/setenv.txt
share/doc/lighttpd/simple-vhost.txt
share/doc/lighttpd/skeleton.txt
share/doc/lighttpd/ssi.txt
share/doc/lighttpd/ssl.txt
share/doc/lighttpd/state.dot
share/doc/lighttpd/state.txt
share/doc/lighttpd/status.txt
share/doc/lighttpd/traffic-shaping.txt
share/doc/lighttpd/trigger_b4_dl.txt
share/doc/lighttpd/userdir.txt
share/doc/lighttpd/webdav.txt
share/examples/lighttpd/lighttpd.conf
share/examples/lighttpd/rrdtool-graph.sh
share/examples/lighttpd/spawn-php.sh
share/examples/rc.d/lighttpd
