# x11/wdm/PLIST
@comment $NetBSD: PLIST,v 1.5 2005/05/02 20:34:08 reed Exp $
bin/wdm
bin/wdmLogin
man/man1/wdm.1
man/man1/wdmLogin.1
share/examples/rc.d/wdm
share/examples/wdm/GiveConsole
share/examples/wdm/TakeConsole
share/examples/wdm/Xaccess
share/examples/wdm/Xclients
share/examples/wdm/Xclients.in
share/examples/wdm/Xresources
share/examples/wdm/Xservers
share/examples/wdm/Xservers.fs
share/examples/wdm/Xservers.ws
share/examples/wdm/Xsession
share/examples/wdm/Xsession.XFree86
share/examples/wdm/Xsession.orig
share/examples/wdm/Xsetup_0
share/examples/wdm/wdm-config
share/examples/wdm/wdm-config.in
share/examples/wdm/wdmReconfig
${PKGLOCALEDIR}/locale/de/LC_MESSAGES/wdm.mo
${PKGLOCALEDIR}/locale/ja/LC_MESSAGES/wdm.mo
${PKGLOCALEDIR}/locale/ru/LC_MESSAGES/wdm.mo
share/pixmaps/wdm/Eagle.xpm
share/pixmaps/wdm/Earth.jpg
share/pixmaps/wdm/Earth200x130.xpm
share/pixmaps/wdm/GnuLinux-logo.jpg
share/pixmaps/wdm/Stars.jpg
share/pixmaps/wdm/StarsLogo.jpg
share/pixmaps/wdm/gnuLogo.xpm
@dirrm share/pixmaps/wdm
@dirrm share/examples/wdm
