# net/chrony/PLIST
@comment $NetBSD: PLIST,v 1.3 2005/05/02 20:34:01 reed Exp $
bin/chronyc
share/doc/chrony/COPYING
share/doc/chrony/README
share/doc/chrony/chrony.txt
share/examples/chrony/chrony.conf.example
share/examples/chrony/chrony.keys.example
man/man1/chrony.1
man/man1/chronyc.1
man/man5/chrony.conf.5
man/man8/chronyd.8
sbin/chronyd
share/examples/rc.d/chronyd
@dirrm share/doc/chrony
@dirrm share/examples/chrony
