# multimedia/ggrab/PLIST
@comment $NetBSD: PLIST,v 1.2 2009/06/14 18:08:04 joerg Exp $
bin/ggrab
bin/sserver
share/doc/ggrab/README
share/examples/rc.d/sserver
