# net/pureftpd/files/pure_ftpd
#!/bin/sh
#
#       $NetBSD: pure_ftpd,v 1.1 2002/05/08 00:35:14 hubertf Exp $
#
# PROVIDE: pure-ftpd
# REQUIRE: DAEMON

. /etc/rc.subr

name="pure_ftpd"
help_name="pure-ftpd"
rcvar=$rcname
command="@PREFIX@/sbin/${help_name}"
command_args="-B"
pidfile="/var/run/${help_name}.pid"

load_rc_config $name
run_rc_command "$1"
