# audio/nas/PLIST.server
@comment $NetBSD: PLIST.server,v 1.2 2008/10/11 21:00:53 markd Exp $
bin/nasd
${IMAKE_MAN_DIR}/nasd.${IMAKE_MANNEWSUFFIX}
${IMAKE_FILEMAN_DIR}/nasd.conf.${IMAKE_FILEMAN_SUFFIX}
share/examples/nas/nasd.conf.eg
share/examples/rc.d/nasd
@dirrm share/examples/nas
