# audio/nas/PLIST.server
@comment $NetBSD: PLIST.server,v 1.3 2009/06/14 17:32:19 joerg Exp $
bin/nasd
${IMAKE_MAN_DIR}/nasd.${IMAKE_MANNEWSUFFIX}
${IMAKE_FILEMAN_DIR}/nasd.conf.${IMAKE_FILEMAN_SUFFIX}
share/examples/nas/nasd.conf.eg
share/examples/rc.d/nasd
