# net/pxe/files/pxe
#!/bin/sh
#
#	$NetBSD: pxe,v 1.2 2001/11/06 05:55:05 cjs Exp $
#

# PROVIDE: pxe
# REQUIRE: DAEMON

. /etc/rc.subr

name="pxe"
command="@@PREFIX@@/sbin/${name}"
required_files="/etc/${name}.conf"

load_rc_config $name
run_rc_command "$1"
