# databases/postgresql/pkg/PLIST
@comment $NetBSD: PLIST,v 1.6 1998/06/22 14:00:16 frueauf Exp $
etc/rc.d/pgsql.sh
@PGUSER@/.profile
@PGUSER@/bin/cleardbdir
@PGUSER@/bin/createdb
@PGUSER@/bin/createuser
@PGUSER@/bin/destroydb
@PGUSER@/bin/destroyuser
@PGUSER@/bin/ecpg
@PGUSER@/bin/initdb
@PGUSER@/bin/initlocation
@PGUSER@/bin/ipcclean
@PGUSER@/bin/pg_dump
@PGUSER@/bin/pg_dumpall
@PGUSER@/bin/pg_id
@PGUSER@/bin/pg_passwd
@PGUSER@/bin/pg_version
@unexec rm -f %D/@PGUSER@/bin/pgtclsh
@unexec rm -f %D/@PGUSER@/bin/pgtksh
@unexec rm -f %D/@PGUSER@/bin/pgaccess
@PGUSER@/bin/postgres
@PGUSER@/bin/postmaster
@PGUSER@/bin/psql
@PGUSER@/data/PG_VERSION
@PGUSER@/data/base/template1/PG_VERSION
@PGUSER@/data/base/template1/pg_aggregate
@PGUSER@/data/base/template1/pg_am
@PGUSER@/data/base/template1/pg_amop
@PGUSER@/data/base/template1/pg_amproc
@PGUSER@/data/base/template1/pg_attrdef
@PGUSER@/data/base/template1/pg_attribute
@PGUSER@/data/base/template1/pg_class
@PGUSER@/data/base/template1/pg_index
@PGUSER@/data/base/template1/pg_inheritproc
@PGUSER@/data/base/template1/pg_inherits
@PGUSER@/data/base/template1/pg_internal.init
@PGUSER@/data/base/template1/pg_ipl
@PGUSER@/data/base/template1/pg_language
@PGUSER@/data/base/template1/pg_listener
@PGUSER@/data/base/template1/pg_opclass
@PGUSER@/data/base/template1/pg_operator
@PGUSER@/data/base/template1/pg_parg
@PGUSER@/data/base/template1/pg_proc
@PGUSER@/data/base/template1/pg_relcheck
@PGUSER@/data/base/template1/pg_rewrite
@PGUSER@/data/base/template1/pg_statistic
@PGUSER@/data/base/template1/pg_trigger
@PGUSER@/data/base/template1/pg_type
@PGUSER@/data/base/template1/pg_version
@PGUSER@/data/base/template1/pg_description
@PGUSER@/data/base/template1/pg_attribute_relid_attnam_index
@PGUSER@/data/base/template1/pg_proc_oid_index
@PGUSER@/data/base/template1/pg_attribute_relid_attnum_index
@PGUSER@/data/base/template1/pg_attribute_attrelid_index
@PGUSER@/data/base/template1/pg_proc_proname_index
@PGUSER@/data/base/template1/pg_proc_prosrc_index
@PGUSER@/data/base/template1/pg_type_oid_index
@PGUSER@/data/base/template1/pg_type_typname_index
@PGUSER@/data/base/template1/pg_class_oid_index
@PGUSER@/data/base/template1/pg_class_relname_index
@PGUSER@/data/base/template1/pg_attrdef_adrelid_index
@PGUSER@/data/base/template1/pg_relcheck_rcrelid_index
@PGUSER@/data/base/template1/pg_trigger_tgrelid_index
@PGUSER@/data/base/template1/pg_description_objoid_index
@PGUSER@/data/base/template1/pg_user
@PGUSER@/data/pg_shadow
@PGUSER@/data/pg_pwd
@PGUSER@/data/pg_database
@PGUSER@/data/pg_geqo.sample
@PGUSER@/data/pg_group
@PGUSER@/data/pg_hba.conf
@PGUSER@/data/pg_log
@PGUSER@/data/pg_variable
@PGUSER@/include/access/attnum.h
@PGUSER@/include/c.h
@PGUSER@/include/config.h
@PGUSER@/include/commands/trigger.h
@PGUSER@/include/executor/spi.h
@PGUSER@/include/fmgr.h
@PGUSER@/include/lib/dllist.h
@unexec rm -f %D/@PGUSER@/include/libpgtcl.h
@PGUSER@/include/libpq-fe.h
@PGUSER@/include/libpq/libpq-fs.h
@PGUSER@/include/libpq/pqcomm.h
@PGUSER@/include/os.h
@PGUSER@/include/postgres.h
@PGUSER@/include/postgres_ext.h
@PGUSER@/include/utils/elog.h
@PGUSER@/include/utils/geo_decls.h
@PGUSER@/include/utils/palloc.h
@PGUSER@/include/ecpglib.h
@PGUSER@/include/ecpgtype.h
@PGUSER@/include/sqlca.h
@PGUSER@/lib/global1.bki.source
@unexec rm -f %D/@PGUSER@/lib/libpgtcl.a
@unexec rm -f %D/@PGUSER@/lib/libpgtcl.so
@unexec rm -f %D/@PGUSER@/lib/libpgtcl.so.1.0
@PGUSER@/lib/libpq.a
@PGUSER@/lib/libpq.so
@PGUSER@/lib/local1_template1.bki.source
@PGUSER@/lib/pg_geqo.sample
@PGUSER@/lib/pg_hba.conf.sample
@PGUSER@/lib/global1.description
@PGUSER@/lib/local1_template1.description
@PGUSER@/lib/libpq.so.1.1
@PGUSER@/lib/libecpg.so.1.1
@PGUSER@/lib/libecpg.so
@PGUSER@/lib/libecpg.a
@PGUSER@/man/man1/cleardbdir.1.gz
@PGUSER@/man/man1/createdb.1.gz
@PGUSER@/man/man1/createuser.1.gz
@PGUSER@/man/man1/destroydb.1.gz
@PGUSER@/man/man1/destroyuser.1.gz
@PGUSER@/man/man1/initdb.1.gz
@PGUSER@/man/man1/ipcclean.1.gz
@PGUSER@/man/man1/pg_dump.1.gz
@PGUSER@/man/man1/pg_dumpall.1.gz
@PGUSER@/man/man1/pg_passwd.1.gz
@PGUSER@/man/man1/pgintro.1.gz
@PGUSER@/man/man1/postgres.1.gz
@PGUSER@/man/man1/postmaster.1.gz
@PGUSER@/man/man1/psql.1.gz
@PGUSER@/man/man3/catalogs.3.gz
@PGUSER@/man/man3/large_objects.3.gz
@PGUSER@/man/man3/libpq.3.gz
@PGUSER@/man/man3/oracle_compat.3.gz
@PGUSER@/man/man3/pgbuiltin.3.gz
@PGUSER@/man/man5/bki.5.gz
@PGUSER@/man/man5/page.5.gz
@PGUSER@/man/man5/pg_hba.conf.5.gz
@PGUSER@/man/manl/abort.l.gz
@PGUSER@/man/manl/alter_table.l.gz
@PGUSER@/man/manl/begin.l.gz
@PGUSER@/man/manl/close.l.gz
@PGUSER@/man/manl/cluster.l.gz
@PGUSER@/man/manl/commit.l.gz
@PGUSER@/man/manl/copy.l.gz
@PGUSER@/man/manl/create_aggregate.l.gz
@PGUSER@/man/manl/create_database.l.gz
@PGUSER@/man/manl/create_function.l.gz
@PGUSER@/man/manl/create_index.l.gz
@PGUSER@/man/manl/create_operator.l.gz
@PGUSER@/man/manl/create_rule.l.gz
@PGUSER@/man/manl/create_sequence.l.gz
@PGUSER@/man/manl/create_table.l.gz
@PGUSER@/man/manl/create_trigger.l.gz
@PGUSER@/man/manl/create_type.l.gz
@PGUSER@/man/manl/create_version.l.gz
@PGUSER@/man/manl/create_view.l.gz
@PGUSER@/man/manl/declare.l.gz
@PGUSER@/man/manl/delete.l.gz
@PGUSER@/man/manl/drop.l.gz
@PGUSER@/man/manl/drop_aggregate.l.gz
@PGUSER@/man/manl/drop_database.l.gz
@PGUSER@/man/manl/drop_function.l.gz
@PGUSER@/man/manl/drop_index.l.gz
@PGUSER@/man/manl/drop_operator.l.gz
@PGUSER@/man/manl/drop_rule.l.gz
@PGUSER@/man/manl/drop_sequence.l.gz
@PGUSER@/man/manl/drop_table.l.gz
@PGUSER@/man/manl/drop_trigger.l.gz
@PGUSER@/man/manl/drop_type.l.gz
@PGUSER@/man/manl/drop_view.l.gz
@PGUSER@/man/manl/end.l.gz
@PGUSER@/man/manl/explain.l.gz
@PGUSER@/man/manl/fetch.l.gz
@PGUSER@/man/manl/grant.l.gz
@PGUSER@/man/manl/insert.l.gz
@PGUSER@/man/manl/listen.l.gz
@PGUSER@/man/manl/load.l.gz
@PGUSER@/man/manl/move.l.gz
@PGUSER@/man/manl/notify.l.gz
@PGUSER@/man/manl/reset.l.gz
@PGUSER@/man/manl/revoke.l.gz
@PGUSER@/man/manl/rollback.l.gz
@PGUSER@/man/manl/select.l.gz
@PGUSER@/man/manl/set.l.gz
@PGUSER@/man/manl/show.l.gz
@PGUSER@/man/manl/sql.l.gz
@PGUSER@/man/manl/update.l.gz
@PGUSER@/man/manl/vacuum.l.gz
@PGUSER@/man/man1/initlocation.1
@PGUSER@/man/man1/ecpg.1
@PGUSER@/man/manl/lock.l
@PGUSER@/man/manl/drop_user.l
@PGUSER@/man/manl/drop_language.l
@PGUSER@/man/manl/create_user.l
@PGUSER@/man/manl/create_language.l
@PGUSER@/man/manl/alter_user.l
@PGUSER@/post-install-notes
share/doc/pgsql/FAQ
share/doc/pgsql/FAQ_Irix
share/doc/pgsql/FAQ_Linux
share/doc/pgsql/README.GEQO
share/doc/pgsql/README.fsync
share/doc/pgsql/README.support
share/doc/pgsql/TODO
share/doc/pgsql/TODO.GEQO
share/doc/pgsql/bug.template
share/doc/pgsql/userguide.ps
share/doc/pgsql/FAQ_FreeBSD
share/doc/pgsql/Makefile
share/doc/pgsql/README.mb
share/doc/pgsql/README.mb.jp
share/doc/pgsql/admin.ps.gz
share/doc/pgsql/admin.tar.gz
share/doc/pgsql/postgres.tar.gz
share/doc/pgsql/programmer.ps.gz
share/doc/pgsql/programmer.tar.gz
share/doc/pgsql/tutorial.ps.gz
share/doc/pgsql/tutorial.tar.gz
share/doc/pgsql/user.ps.gz
share/doc/pgsql/user.tar.gz
share/doc/pgsql/src/Makefile
share/doc/pgsql/src/graphics/
share/doc/pgsql/src/graphics/catalogs.gif
share/doc/pgsql/src/graphics/clientserver.ag
share/doc/pgsql/src/graphics/clientserver.gif
share/doc/pgsql/src/graphics/connections.ag
share/doc/pgsql/src/graphics/connections.gif
share/doc/pgsql/src/graphics/layout.ag
share/doc/pgsql/src/graphics/layout.gif
share/doc/pgsql/src/sgml/Makefile
share/doc/pgsql/src/sgml/admin.sgml
share/doc/pgsql/src/sgml/advanced.sgml
share/doc/pgsql/src/sgml/arch-dev.sgml
share/doc/pgsql/src/sgml/arch-pg.sgml
share/doc/pgsql/src/sgml/arch.sgml
share/doc/pgsql/src/sgml/array.sgml
share/doc/pgsql/src/sgml/biblio.sgml
share/doc/pgsql/src/sgml/compiler.sgml
share/doc/pgsql/src/sgml/contacts.sgml
share/doc/pgsql/src/sgml/datatype.sgml
share/doc/pgsql/src/sgml/dfunc.sgml
share/doc/pgsql/src/sgml/docguide.sgml
share/doc/pgsql/src/sgml/ecpg.sgml
share/doc/pgsql/src/sgml/environ.sgml
share/doc/pgsql/src/sgml/extend.sgml
share/doc/pgsql/src/sgml/func-ref.sgml
share/doc/pgsql/src/sgml/geqo.sgml
share/doc/pgsql/src/sgml/gist.sgml
share/doc/pgsql/src/sgml/inherit.sgml
share/doc/pgsql/src/sgml/install.sgml
share/doc/pgsql/src/sgml/intro-pg.sgml
share/doc/pgsql/src/sgml/intro.sgml
share/doc/pgsql/src/sgml/jdbc.sgml
share/doc/pgsql/src/sgml/libpgtcl.sgml
share/doc/pgsql/src/sgml/libpq.sgml
share/doc/pgsql/src/sgml/lobj.sgml
share/doc/pgsql/src/sgml/manage.sgml
share/doc/pgsql/src/sgml/odbc.sgml
share/doc/pgsql/src/sgml/pgaccess.sgml
share/doc/pgsql/src/sgml/ports.sgml
share/doc/pgsql/src/sgml/postgres.sgml
share/doc/pgsql/src/sgml/programmer.sgml
share/doc/pgsql/src/sgml/protocol.sgml
share/doc/pgsql/src/sgml/psql.sgml
share/doc/pgsql/src/sgml/query-ug.sgml
share/doc/pgsql/src/sgml/query.sgml
share/doc/pgsql/src/sgml/recovery.sgml
share/doc/pgsql/src/sgml/regress.sgml
share/doc/pgsql/src/sgml/release.sgml
share/doc/pgsql/src/sgml/rules.sgml
share/doc/pgsql/src/sgml/spi.sgml
share/doc/pgsql/src/sgml/start-ag.sgml
share/doc/pgsql/src/sgml/start.sgml
share/doc/pgsql/src/sgml/storage.sgml
share/doc/pgsql/src/sgml/trigger.sgml
share/doc/pgsql/src/sgml/tutorial.sgml
share/doc/pgsql/src/sgml/user.sgml
share/doc/pgsql/src/sgml/xaggr.sgml
share/doc/pgsql/src/sgml/xfunc.sgml
share/doc/pgsql/src/sgml/xindex.sgml
share/doc/pgsql/src/sgml/xoper.sgml
share/doc/pgsql/src/sgml/xtypes.sgml
@unexec /sbin/ldconfig || /usr/bin/true
@dirrm share/doc/pgsql/src/graphics
@dirrm share/doc/pgsql/src/sgml
@dirrm share/doc/pgsql/src
@dirrm share/doc/pgsql
@dirrm @PGUSER@/bin
@dirrm @PGUSER@/data/base/template1
@dirrm @PGUSER@/data/base
@dirrm @PGUSER@/data
@dirrm @PGUSER@/include/access
@dirrm @PGUSER@/include/commands
@dirrm @PGUSER@/include/executor
@dirrm @PGUSER@/include/include
@dirrm @PGUSER@/include/lib
@dirrm @PGUSER@/include/libpq
@dirrm @PGUSER@/include/port/bsd
@dirrm @PGUSER@/include/port
@dirrm @PGUSER@/include/utils
@dirrm @PGUSER@/include
@dirrm @PGUSER@/lib
@dirrm @PGUSER@/man/man1
@dirrm @PGUSER@/man/man3
@dirrm @PGUSER@/man/man5
@dirrm @PGUSER@/man/manl
@dirrm @PGUSER@/man
@dirrm @PGUSER@
