# databases/postgresql/pkg/PLIST
@comment $NetBSD: PLIST,v 1.8 1999/01/04 13:05:02 frueauf Exp $
etc/rc.d/pgsql.sh
share/doc/pgsql/src/graphics/catalogs.gif
share/doc/pgsql/src/graphics/catalogs.ps
share/doc/pgsql/src/graphics/clientserver.ag
share/doc/pgsql/src/graphics/clientserver.gif
share/doc/pgsql/src/graphics/connections.ag
share/doc/pgsql/src/graphics/connections.gif
share/doc/pgsql/src/graphics/layout.ag
share/doc/pgsql/src/graphics/layout.gif
share/doc/pgsql/src/sgml/ref/abort.sgml
share/doc/pgsql/src/sgml/ref/allfiles.sgml
share/doc/pgsql/src/sgml/ref/alter_table.sgml
share/doc/pgsql/src/sgml/ref/alter_user.sgml
share/doc/pgsql/src/sgml/ref/begin.sgml
share/doc/pgsql/src/sgml/ref/close.sgml
share/doc/pgsql/src/sgml/ref/cluster.sgml
share/doc/pgsql/src/sgml/ref/commands.sgml
share/doc/pgsql/src/sgml/ref/commit.sgml
share/doc/pgsql/src/sgml/ref/copy.sgml
share/doc/pgsql/src/sgml/ref/create_aggregate.sgml
share/doc/pgsql/src/sgml/ref/create_database.sgml
share/doc/pgsql/src/sgml/ref/create_function.sgml
share/doc/pgsql/src/sgml/ref/create_index.sgml
share/doc/pgsql/src/sgml/ref/create_language.sgml
share/doc/pgsql/src/sgml/ref/create_operator.sgml
share/doc/pgsql/src/sgml/ref/create_rule.sgml
share/doc/pgsql/src/sgml/ref/create_sequence.sgml
share/doc/pgsql/src/sgml/ref/create_table.sgml
share/doc/pgsql/src/sgml/ref/create_trigger.sgml
share/doc/pgsql/src/sgml/ref/create_type.sgml
share/doc/pgsql/src/sgml/ref/create_user.sgml
share/doc/pgsql/src/sgml/ref/create_view.sgml
share/doc/pgsql/src/sgml/ref/createdb.sgml
share/doc/pgsql/src/sgml/ref/createuser.sgml
share/doc/pgsql/src/sgml/ref/current_date.sgml
share/doc/pgsql/src/sgml/ref/current_time.sgml
share/doc/pgsql/src/sgml/ref/current_timestamp.sgml
share/doc/pgsql/src/sgml/ref/current_user.sgml
share/doc/pgsql/src/sgml/ref/declare.sgml
share/doc/pgsql/src/sgml/ref/delete.sgml
share/doc/pgsql/src/sgml/ref/destroydb.sgml
share/doc/pgsql/src/sgml/ref/destroyuser.sgml
share/doc/pgsql/src/sgml/ref/drop_aggregate.sgml
share/doc/pgsql/src/sgml/ref/drop_database.sgml
share/doc/pgsql/src/sgml/ref/drop_function.sgml
share/doc/pgsql/src/sgml/ref/drop_index.sgml
share/doc/pgsql/src/sgml/ref/drop_language.sgml
share/doc/pgsql/src/sgml/ref/drop_operator.sgml
share/doc/pgsql/src/sgml/ref/drop_rule.sgml
share/doc/pgsql/src/sgml/ref/drop_sequence.sgml
share/doc/pgsql/src/sgml/ref/drop_table.sgml
share/doc/pgsql/src/sgml/ref/drop_trigger.sgml
share/doc/pgsql/src/sgml/ref/drop_type.sgml
share/doc/pgsql/src/sgml/ref/drop_user.sgml
share/doc/pgsql/src/sgml/ref/drop_view.sgml
share/doc/pgsql/src/sgml/ref/explain.sgml
share/doc/pgsql/src/sgml/ref/fetch.sgml
share/doc/pgsql/src/sgml/ref/grant.sgml
share/doc/pgsql/src/sgml/ref/initdb.sgml
share/doc/pgsql/src/sgml/ref/initlocation.sgml
share/doc/pgsql/src/sgml/ref/insert.sgml
share/doc/pgsql/src/sgml/ref/listen.sgml
share/doc/pgsql/src/sgml/ref/load.sgml
share/doc/pgsql/src/sgml/ref/lock.sgml
share/doc/pgsql/src/sgml/ref/move.sgml
share/doc/pgsql/src/sgml/ref/notify.sgml
share/doc/pgsql/src/sgml/ref/pg_dump.sgml
share/doc/pgsql/src/sgml/ref/pg_dumpall.sgml
share/doc/pgsql/src/sgml/ref/pg_upgrade.sgml
share/doc/pgsql/src/sgml/ref/psql-ref.sgml
share/doc/pgsql/src/sgml/ref/reset.sgml
share/doc/pgsql/src/sgml/ref/revoke.sgml
share/doc/pgsql/src/sgml/ref/rollback.sgml
share/doc/pgsql/src/sgml/ref/select.sgml
share/doc/pgsql/src/sgml/ref/set.sgml
share/doc/pgsql/src/sgml/ref/show.sgml
share/doc/pgsql/src/sgml/ref/unlisten.sgml
share/doc/pgsql/src/sgml/ref/update.sgml
share/doc/pgsql/src/sgml/ref/vacuum.sgml
share/doc/pgsql/src/sgml/about.sgml
share/doc/pgsql/src/sgml/admin.sgml
share/doc/pgsql/src/sgml/advanced.sgml
share/doc/pgsql/src/sgml/arch-dev.sgml
share/doc/pgsql/src/sgml/arch-pg.sgml
share/doc/pgsql/src/sgml/arch.sgml
share/doc/pgsql/src/sgml/array.sgml
share/doc/pgsql/src/sgml/biblio.sgml
share/doc/pgsql/src/sgml/bki.sgml
share/doc/pgsql/src/sgml/compiler.sgml
share/doc/pgsql/src/sgml/config.sgml
share/doc/pgsql/src/sgml/contacts.sgml
share/doc/pgsql/src/sgml/current.sgml
share/doc/pgsql/src/sgml/datatype.sgml
share/doc/pgsql/src/sgml/dfunc.sgml
share/doc/pgsql/src/sgml/docguide.sgml
share/doc/pgsql/src/sgml/ecpg.sgml
share/doc/pgsql/src/sgml/environ.sgml
share/doc/pgsql/src/sgml/extend.sgml
share/doc/pgsql/src/sgml/func-ref.sgml
share/doc/pgsql/src/sgml/func.sgml
share/doc/pgsql/src/sgml/geqo.sgml
share/doc/pgsql/src/sgml/gist.sgml
share/doc/pgsql/src/sgml/history.sgml
share/doc/pgsql/src/sgml/indices.sgml
share/doc/pgsql/src/sgml/info.sgml
share/doc/pgsql/src/sgml/inherit.sgml
share/doc/pgsql/src/sgml/install.sgml
share/doc/pgsql/src/sgml/installation.sgml
share/doc/pgsql/src/sgml/intro-ag.sgml
share/doc/pgsql/src/sgml/intro-pg.sgml
share/doc/pgsql/src/sgml/intro.sgml
share/doc/pgsql/src/sgml/jdbc.sgml
share/doc/pgsql/src/sgml/keys.sgml
share/doc/pgsql/src/sgml/legal.sgml
share/doc/pgsql/src/sgml/libpgtcl.sgml
share/doc/pgsql/src/sgml/libpq++.sgml
share/doc/pgsql/src/sgml/libpq.sgml
share/doc/pgsql/src/sgml/lobj.sgml
share/doc/pgsql/src/sgml/manage.sgml
share/doc/pgsql/src/sgml/notation.sgml
share/doc/pgsql/src/sgml/odbc.sgml
share/doc/pgsql/src/sgml/oper.sgml
share/doc/pgsql/src/sgml/page.sgml
share/doc/pgsql/src/sgml/pg_options.sgml
share/doc/pgsql/src/sgml/pgaccess.sgml
share/doc/pgsql/src/sgml/ports.sgml
share/doc/pgsql/src/sgml/postgres.sgml
share/doc/pgsql/src/sgml/programmer.sgml
share/doc/pgsql/src/sgml/protocol.sgml
share/doc/pgsql/src/sgml/psql.sgml
share/doc/pgsql/src/sgml/query-ug.sgml
share/doc/pgsql/src/sgml/query.sgml
share/doc/pgsql/src/sgml/recovery.sgml
share/doc/pgsql/src/sgml/reference.sgml
share/doc/pgsql/src/sgml/regress.sgml
share/doc/pgsql/src/sgml/release.sgml
share/doc/pgsql/src/sgml/rules.sgml
share/doc/pgsql/src/sgml/runtime.sgml
share/doc/pgsql/src/sgml/security.sgml
share/doc/pgsql/src/sgml/signals.sgml
share/doc/pgsql/src/sgml/spi.sgml
share/doc/pgsql/src/sgml/start-ag.sgml
share/doc/pgsql/src/sgml/start.sgml
share/doc/pgsql/src/sgml/storage.sgml
share/doc/pgsql/src/sgml/syntax.sgml
share/doc/pgsql/src/sgml/trigger.sgml
share/doc/pgsql/src/sgml/tutorial.sgml
share/doc/pgsql/src/sgml/typeconv.sgml
share/doc/pgsql/src/sgml/user.sgml
share/doc/pgsql/src/sgml/xaggr.sgml
share/doc/pgsql/src/sgml/xfunc.sgml
share/doc/pgsql/src/sgml/xindex.sgml
share/doc/pgsql/src/sgml/xoper.sgml
share/doc/pgsql/src/sgml/xplang.sgml
share/doc/pgsql/src/sgml/xtypes.sgml
share/doc/pgsql/src/sgml/y2k.sgml
share/doc/pgsql/src/sgml/Makefile
share/doc/pgsql/src/Makefile
share/doc/pgsql/tutorial.tar.gz
share/doc/pgsql/tutorial.ps.gz
share/doc/pgsql/programmer.tar.gz
share/doc/pgsql/programmer.ps.gz
share/doc/pgsql/postgres.tar.gz
share/doc/pgsql/bug.template
share/doc/pgsql/admin.tar.gz
share/doc/pgsql/admin.ps.gz
share/doc/pgsql/TODO
share/doc/pgsql/README.mb.jp
share/doc/pgsql/user.tar.gz
share/doc/pgsql/README.locale
share/doc/pgsql/README.inet
share/doc/pgsql/README.fsync
share/doc/pgsql/Makefile
share/doc/pgsql/FAQ_Linux
share/doc/pgsql/FAQ_Irix
share/doc/pgsql/FAQ_HPUX
share/doc/pgsql/FAQ_FreeBSD
share/doc/pgsql/FAQ_DEV
share/doc/pgsql/FAQ_CVS
share/doc/pgsql/FAQ
share/doc/pgsql/README.mb
share/doc/pgsql/user.ps.gz
@PGUSER@/lib/global1.bki.source
@PGUSER@/lib/global1.description
@PGUSER@/lib/local1_template1.bki.source
@PGUSER@/lib/local1_template1.description
@PGUSER@/lib/pg_hba.conf.sample
@PGUSER@/lib/pg_geqo.sample
@PGUSER@/lib/libpq.a
@PGUSER@/lib/libpq.so.2.0
@PGUSER@/lib/libpq.so.2
@PGUSER@/lib/libpq.so
@PGUSER@/lib/libecpg.a
@PGUSER@/lib/libecpg.so.2.6.2
@PGUSER@/lib/libecpg.so.2
@PGUSER@/lib/libecpg.so
@PGUSER@/lib/libpq++.a
@PGUSER@/lib/libpq++.so.2.0
@PGUSER@/lib/libpq++.so.2
@PGUSER@/lib/libpq++.so
@PGUSER@/lib/libpgtcl.a
@PGUSER@/lib/libpgtcl.so.2.0
@PGUSER@/lib/libpgtcl.so.2
@PGUSER@/lib/libpgtcl.so
@PGUSER@/lib/plpgsql.so
@PGUSER@/lib/pltcl.so
@PGUSER@/bin/postgres
@PGUSER@/bin/postmaster
@PGUSER@/bin/ecpg
@PGUSER@/bin/pg_id
@PGUSER@/bin/pg_version
@PGUSER@/bin/psql
@PGUSER@/bin/pg_dump
@PGUSER@/bin/pg_dumpall
@PGUSER@/bin/pg_upgrade
@PGUSER@/bin/pg_passwd
@PGUSER@/bin/cleardbdir
@PGUSER@/bin/createdb
@PGUSER@/bin/createuser
@PGUSER@/bin/destroydb
@PGUSER@/bin/destroyuser
@PGUSER@/bin/initdb
@PGUSER@/bin/initlocation
@PGUSER@/bin/ipcclean
@PGUSER@/bin/pgaccess
@PGUSER@/bin/pgtclsh
@PGUSER@/bin/pgtksh
@PGUSER@/man/man1/cleardbdir.1
@PGUSER@/man/man1/createdb.1
@PGUSER@/man/man1/createuser.1
@PGUSER@/man/man1/destroydb.1
@PGUSER@/man/man1/destroyuser.1
@PGUSER@/man/man1/ecpg.1
@PGUSER@/man/man1/initdb.1
@PGUSER@/man/man1/initlocation.1
@PGUSER@/man/man1/ipcclean.1
@PGUSER@/man/man1/pg_dump.1
@PGUSER@/man/man1/pg_dumpall.1
@PGUSER@/man/man1/pg_passwd.1
@PGUSER@/man/man1/pg_upgrade.1
@PGUSER@/man/man1/postgres.1
@PGUSER@/man/man1/postmaster.1
@PGUSER@/man/man1/psql.1
@PGUSER@/man/man3/catalogs.3
@PGUSER@/man/man3/libpq.3
@PGUSER@/man/man5/pg_hba.conf.5
@PGUSER@/man/manl/abort.l
@PGUSER@/man/manl/alter_table.l
@PGUSER@/man/manl/alter_user.l
@PGUSER@/man/manl/begin.l
@PGUSER@/man/manl/close.l
@PGUSER@/man/manl/cluster.l
@PGUSER@/man/manl/commit.l
@PGUSER@/man/manl/copy.l
@PGUSER@/man/manl/create_aggregate.l
@PGUSER@/man/manl/create_database.l
@PGUSER@/man/manl/create_function.l
@PGUSER@/man/manl/create_index.l
@PGUSER@/man/manl/create_language.l
@PGUSER@/man/manl/create_operator.l
@PGUSER@/man/manl/create_rule.l
@PGUSER@/man/manl/create_sequence.l
@PGUSER@/man/manl/create_table.l
@PGUSER@/man/manl/create_trigger.l
@PGUSER@/man/manl/create_type.l
@PGUSER@/man/manl/create_user.l
@PGUSER@/man/manl/create_version.l
@PGUSER@/man/manl/create_view.l
@PGUSER@/man/manl/declare.l
@PGUSER@/man/manl/delete.l
@PGUSER@/man/manl/drop.l
@PGUSER@/man/manl/drop_aggregate.l
@PGUSER@/man/manl/drop_database.l
@PGUSER@/man/manl/drop_function.l
@PGUSER@/man/manl/drop_index.l
@PGUSER@/man/manl/drop_language.l
@PGUSER@/man/manl/drop_operator.l
@PGUSER@/man/manl/drop_rule.l
@PGUSER@/man/manl/drop_sequence.l
@PGUSER@/man/manl/drop_table.l
@PGUSER@/man/manl/drop_trigger.l
@PGUSER@/man/manl/drop_type.l
@PGUSER@/man/manl/drop_user.l
@PGUSER@/man/manl/drop_view.l
@PGUSER@/man/manl/end.l
@PGUSER@/man/manl/explain.l
@PGUSER@/man/manl/fetch.l
@PGUSER@/man/manl/grant.l
@PGUSER@/man/manl/insert.l
@PGUSER@/man/manl/listen.l
@PGUSER@/man/manl/load.l
@PGUSER@/man/manl/lock.l
@PGUSER@/man/manl/move.l
@PGUSER@/man/manl/notify.l
@PGUSER@/man/manl/reset.l
@PGUSER@/man/manl/revoke.l
@PGUSER@/man/manl/rollback.l
@PGUSER@/man/manl/select.l
@PGUSER@/man/manl/set.l
@PGUSER@/man/manl/show.l
@PGUSER@/man/manl/sql.l
@PGUSER@/man/manl/update.l
@PGUSER@/man/manl/vacuum.l
@PGUSER@/data/base/template1/pg_proc
@PGUSER@/data/base/template1/pg_type
@PGUSER@/data/base/template1/pg_attribute
@PGUSER@/data/base/template1/pg_class
@PGUSER@/data/base/template1/pg_inherits
@PGUSER@/data/base/template1/pg_index
@PGUSER@/data/base/template1/pg_statistic
@PGUSER@/data/base/template1/pg_operator
@PGUSER@/data/base/template1/pg_opclass
@PGUSER@/data/base/template1/pg_am
@PGUSER@/data/base/template1/pg_amop
@PGUSER@/data/base/template1/pg_amproc
@PGUSER@/data/base/template1/pg_language
@PGUSER@/data/base/template1/pg_parg
@PGUSER@/data/base/template1/pg_aggregate
@PGUSER@/data/base/template1/pg_ipl
@PGUSER@/data/base/template1/pg_inheritproc
@PGUSER@/data/base/template1/pg_rewrite
@PGUSER@/data/base/template1/pg_listener
@PGUSER@/data/base/template1/pg_description
@PGUSER@/data/base/template1/pg_attribute_relid_attnam_index
@PGUSER@/data/base/template1/pg_attribute_relid_attnum_index
@PGUSER@/data/base/template1/pg_attribute_attrelid_index
@PGUSER@/data/base/template1/pg_proc_oid_index
@PGUSER@/data/base/template1/pg_proc_proname_narg_type_index
@PGUSER@/data/base/template1/pg_proc_prosrc_index
@PGUSER@/data/base/template1/pg_type_oid_index
@PGUSER@/data/base/template1/pg_type_typname_index
@PGUSER@/data/base/template1/pg_class_oid_index
@PGUSER@/data/base/template1/pg_class_relname_index
@PGUSER@/data/base/template1/pg_attrdef
@PGUSER@/data/base/template1/pg_attrdef_adrelid_index
@PGUSER@/data/base/template1/pg_relcheck
@PGUSER@/data/base/template1/pg_relcheck_rcrelid_index
@PGUSER@/data/base/template1/pg_trigger
@PGUSER@/data/base/template1/pg_trigger_tgrelid_index
@PGUSER@/data/base/template1/pg_description_objoid_index
@PGUSER@/data/base/template1/PG_VERSION
@PGUSER@/data/base/template1/pg_internal.init
@PGUSER@/data/base/template1/pg_user
@PGUSER@/data/base/template1/pg_rules
@PGUSER@/data/base/template1/pg_views
@PGUSER@/data/base/template1/pg_tables
@PGUSER@/data/base/template1/pg_indexes
@PGUSER@/data/pg_variable
@PGUSER@/data/pg_database
@PGUSER@/data/pg_shadow
@PGUSER@/data/pg_group
@PGUSER@/data/pg_log
@PGUSER@/data/PG_VERSION
@PGUSER@/data/pg_hba.conf
@PGUSER@/data/pg_geqo.sample
@PGUSER@/data/pg_pwd
@PGUSER@/include/fmgr.h
@PGUSER@/include/port/bsd
@PGUSER@/include/lib/dllist.h
@PGUSER@/include/libpq/pqcomm.h
@PGUSER@/include/libpq/libpq-fs.h
@PGUSER@/include/utils/geo_decls.h
@PGUSER@/include/utils/elog.h
@PGUSER@/include/utils/palloc.h
@PGUSER@/include/access/attnum.h
@PGUSER@/include/executor/spi.h
@PGUSER@/include/commands/trigger.h
@PGUSER@/include/os.h
@PGUSER@/include/config.h
@PGUSER@/include/c.h
@PGUSER@/include/postgres.h
@PGUSER@/include/postgres_ext.h
@PGUSER@/include/libpq-fe.h
@PGUSER@/include/libpq-int.h
@PGUSER@/include/ecpgerrno.h
@PGUSER@/include/ecpglib.h
@PGUSER@/include/ecpgtype.h
@PGUSER@/include/sqlca.h
@PGUSER@/include/libpq++/pgenv.h
@PGUSER@/include/libpq++/pgconnection.h
@PGUSER@/include/libpq++/pgdatabase.h
@PGUSER@/include/libpq++/pgtransdb.h
@PGUSER@/include/libpq++/pgcursordb.h
@PGUSER@/include/libpq++/pglobject.h
@PGUSER@/include/libpq++.h
@PGUSER@/include/libpgtcl.h
@PGUSER@/.profile
@PGUSER@/post-install-notes
@dirrm @PGUSER@/data/base/template1
@dirrm @PGUSER@/data/base
@dirrm @PGUSER@/data
@dirrm @PGUSER@/include/libpq++
@dirrm @PGUSER@/include/commands
@dirrm @PGUSER@/include/executor
@dirrm @PGUSER@/include/access
@dirrm @PGUSER@/include/utils
@dirrm @PGUSER@/include/port
@dirrm @PGUSER@/include/libpq
@dirrm @PGUSER@/include/lib
@dirrm @PGUSER@/include
@dirrm @PGUSER@/bin
@dirrm @PGUSER@/lib
@dirrm @PGUSER@/man/man1
@dirrm @PGUSER@/man/man3
@dirrm @PGUSER@/man/man5
@dirrm @PGUSER@/man/manl
@dirrm @PGUSER@/man
@dirrm @PGUSER@
@dirrm share/doc/pgsql/src/sgml/ref
@dirrm share/doc/pgsql/src/sgml
@dirrm share/doc/pgsql/src/graphics
@dirrm share/doc/pgsql/src
@dirrm share/doc/pgsql
