# www/dansguardian/PLIST
@comment $NetBSD: PLIST,v 1.1.1.1 2008/02/22 17:13:29 sborrill Exp $
libexec/cgi-bin/dansguardian.pl
man/man8/dansguardian.8
sbin/dansguardian
share/examples/dansguardian/bannedextensionlist
share/examples/dansguardian/bannediplist
share/examples/dansguardian/bannedmimetypelist
share/examples/dansguardian/bannedphraselist
share/examples/dansguardian/bannedregexpurllist
share/examples/dansguardian/bannedsitelist
share/examples/dansguardian/bannedurllist
share/examples/dansguardian/banneduserlist
share/examples/dansguardian/contentregexplist
share/examples/dansguardian/dansguardian.conf
share/examples/dansguardian/dansguardianf1.conf
share/examples/dansguardian/exceptioniplist
share/examples/dansguardian/exceptionphraselist
share/examples/dansguardian/exceptionsitelist
share/examples/dansguardian/exceptionurllist
share/examples/dansguardian/exceptionuserlist
share/examples/dansguardian/filtergroupslist
share/examples/dansguardian/greysitelist
share/examples/dansguardian/greyurllist
share/examples/dansguardian/languages/ReadMe
share/examples/dansguardian/languages/arspanish/messages
share/examples/dansguardian/languages/arspanish/template.html
share/examples/dansguardian/languages/bulgarian/messages
share/examples/dansguardian/languages/bulgarian/template.html
share/examples/dansguardian/languages/chinesebig5/messages
share/examples/dansguardian/languages/chinesebig5/template.html
share/examples/dansguardian/languages/chinesegb2312/messages
share/examples/dansguardian/languages/chinesegb2312/template.html
share/examples/dansguardian/languages/czech/messages
share/examples/dansguardian/languages/czech/template.html
share/examples/dansguardian/languages/danish/messages
share/examples/dansguardian/languages/danish/template.html
share/examples/dansguardian/languages/dutch/messages
share/examples/dansguardian/languages/dutch/template.html
share/examples/dansguardian/languages/french/messages
share/examples/dansguardian/languages/french/template.html
share/examples/dansguardian/languages/german/messages
share/examples/dansguardian/languages/german/template.html
share/examples/dansguardian/languages/hebrew/messages
share/examples/dansguardian/languages/hebrew/template.html
share/examples/dansguardian/languages/indonesian/messages
share/examples/dansguardian/languages/indonesian/template.html
share/examples/dansguardian/languages/italian/messages
share/examples/dansguardian/languages/italian/template.html
share/examples/dansguardian/languages/lithuanian/messages
share/examples/dansguardian/languages/lithuanian/template.html
share/examples/dansguardian/languages/mxspanish/messages
share/examples/dansguardian/languages/mxspanish/template.html
share/examples/dansguardian/languages/polish/messages
share/examples/dansguardian/languages/polish/template.html
share/examples/dansguardian/languages/portuguese/messages
share/examples/dansguardian/languages/portuguese/messagesaccents
share/examples/dansguardian/languages/portuguese/template.html
share/examples/dansguardian/languages/russian-1251/messages
share/examples/dansguardian/languages/russian-1251/template.html
share/examples/dansguardian/languages/slovak/messages
share/examples/dansguardian/languages/slovak/template.html
share/examples/dansguardian/languages/spanish/messages
share/examples/dansguardian/languages/spanish/template.html
share/examples/dansguardian/languages/swedish/messages
share/examples/dansguardian/languages/swedish/template.html
share/examples/dansguardian/languages/turkish/messages
share/examples/dansguardian/languages/turkish/template.html
share/examples/dansguardian/languages/ukenglish/messages
share/examples/dansguardian/languages/ukenglish/template.html
share/examples/dansguardian/logrotation
share/examples/dansguardian/phraselists/badwords/weighted_dutch
share/examples/dansguardian/phraselists/badwords/weighted_french
share/examples/dansguardian/phraselists/badwords/weighted_german
share/examples/dansguardian/phraselists/badwords/weighted_portuguese
share/examples/dansguardian/phraselists/badwords/weighted_spanish
share/examples/dansguardian/phraselists/chat/weighted
share/examples/dansguardian/phraselists/chat/weighted_italian
share/examples/dansguardian/phraselists/drugadvocacy/weighted
share/examples/dansguardian/phraselists/gambling/banned
share/examples/dansguardian/phraselists/gambling/banned_portuguese
share/examples/dansguardian/phraselists/gambling/weighted
share/examples/dansguardian/phraselists/gambling/weighted_portuguese
share/examples/dansguardian/phraselists/games/weighted
share/examples/dansguardian/phraselists/goodphrases/exception
share/examples/dansguardian/phraselists/goodphrases/exception_email
share/examples/dansguardian/phraselists/goodphrases/weighted_general
share/examples/dansguardian/phraselists/goodphrases/weighted_general_danish
share/examples/dansguardian/phraselists/goodphrases/weighted_general_portuguese
share/examples/dansguardian/phraselists/goodphrases/weighted_news
share/examples/dansguardian/phraselists/googlesearches/banned
share/examples/dansguardian/phraselists/gore/weighted
share/examples/dansguardian/phraselists/gore/weighted_portuguese
share/examples/dansguardian/phraselists/illegaldrugs/banned
share/examples/dansguardian/phraselists/illegaldrugs/banned_portuguese
share/examples/dansguardian/phraselists/illegaldrugs/weighted
share/examples/dansguardian/phraselists/illegaldrugs/weighted_portuguese
share/examples/dansguardian/phraselists/intolerance/banned_portuguese
share/examples/dansguardian/phraselists/intolerance/weighted
share/examples/dansguardian/phraselists/intolerance/weighted_portuguese
share/examples/dansguardian/phraselists/legaldrugs/weighted
share/examples/dansguardian/phraselists/malware/weighted
share/examples/dansguardian/phraselists/news/weighted
share/examples/dansguardian/phraselists/nudism/weighted
share/examples/dansguardian/phraselists/peer2peer/weighted
share/examples/dansguardian/phraselists/personals/weighted
share/examples/dansguardian/phraselists/personals/weighted_portuguese
share/examples/dansguardian/phraselists/pornography/banned
share/examples/dansguardian/phraselists/pornography/banned_portuguese
share/examples/dansguardian/phraselists/pornography/weighted
share/examples/dansguardian/phraselists/pornography/weighted_danish
share/examples/dansguardian/phraselists/pornography/weighted_dutch
share/examples/dansguardian/phraselists/pornography/weighted_french
share/examples/dansguardian/phraselists/pornography/weighted_german
share/examples/dansguardian/phraselists/pornography/weighted_italian
share/examples/dansguardian/phraselists/pornography/weighted_portuguese
share/examples/dansguardian/phraselists/pornography/weighted_spanish
share/examples/dansguardian/phraselists/proxies/weighted
share/examples/dansguardian/phraselists/sport/weighted
share/examples/dansguardian/phraselists/violence/weighted
share/examples/dansguardian/phraselists/violence/weighted_portuguese
share/examples/dansguardian/phraselists/warezhacking/weighted
share/examples/dansguardian/phraselists/weapons/weighted
share/examples/dansguardian/phraselists/weapons/weighted_portuguese
share/examples/dansguardian/phraselists/webmail/weighted
share/examples/dansguardian/pics
share/examples/dansguardian/transparent1x1.gif
share/examples/dansguardian/weightedphraselist
share/examples/rc.d/dansguardian
@dirrm share/examples/dansguardian/phraselists/webmail
@dirrm share/examples/dansguardian/phraselists/weapons
@dirrm share/examples/dansguardian/phraselists/warezhacking
@dirrm share/examples/dansguardian/phraselists/violence
@dirrm share/examples/dansguardian/phraselists/sport
@dirrm share/examples/dansguardian/phraselists/proxies
@dirrm share/examples/dansguardian/phraselists/pornography
@dirrm share/examples/dansguardian/phraselists/personals
@dirrm share/examples/dansguardian/phraselists/peer2peer
@dirrm share/examples/dansguardian/phraselists/nudism
@dirrm share/examples/dansguardian/phraselists/news
@dirrm share/examples/dansguardian/phraselists/malware
@dirrm share/examples/dansguardian/phraselists/legaldrugs
@dirrm share/examples/dansguardian/phraselists/intolerance
@dirrm share/examples/dansguardian/phraselists/illegaldrugs
@dirrm share/examples/dansguardian/phraselists/gore
@dirrm share/examples/dansguardian/phraselists/googlesearches
@dirrm share/examples/dansguardian/phraselists/goodphrases
@dirrm share/examples/dansguardian/phraselists/games
@dirrm share/examples/dansguardian/phraselists/gambling
@dirrm share/examples/dansguardian/phraselists/drugadvocacy
@dirrm share/examples/dansguardian/phraselists/chat
@dirrm share/examples/dansguardian/phraselists/badwords
@dirrm share/examples/dansguardian/phraselists
@dirrm share/examples/dansguardian/languages/ukenglish
@dirrm share/examples/dansguardian/languages/turkish
@dirrm share/examples/dansguardian/languages/swedish
@dirrm share/examples/dansguardian/languages/spanish
@dirrm share/examples/dansguardian/languages/slovak
@dirrm share/examples/dansguardian/languages/russian-1251
@dirrm share/examples/dansguardian/languages/portuguese
@dirrm share/examples/dansguardian/languages/polish
@dirrm share/examples/dansguardian/languages/mxspanish
@dirrm share/examples/dansguardian/languages/lithuanian
@dirrm share/examples/dansguardian/languages/italian
@dirrm share/examples/dansguardian/languages/indonesian
@dirrm share/examples/dansguardian/languages/hebrew
@dirrm share/examples/dansguardian/languages/german
@dirrm share/examples/dansguardian/languages/french
@dirrm share/examples/dansguardian/languages/dutch
@dirrm share/examples/dansguardian/languages/danish
@dirrm share/examples/dansguardian/languages/czech
@dirrm share/examples/dansguardian/languages/chinesegb2312
@dirrm share/examples/dansguardian/languages/chinesebig5
@dirrm share/examples/dansguardian/languages/bulgarian
@dirrm share/examples/dansguardian/languages/arspanish
@dirrm share/examples/dansguardian/languages
@dirrm share/examples/dansguardian
