# net/mldonkey/PLIST
@comment $NetBSD: PLIST,v 1.5 2009/06/14 21:00:04 joerg Exp $
libexec/mldonkey/mlnet
share/examples/rc.d/mldonkey
