# x11/xfstt/PLIST
@comment $NetBSD: PLIST,v 1.7 2005/05/02 20:34:09 reed Exp $
bin/xfstt
man/man1/xfstt.1
share/doc/xfstt/CHANGES
share/doc/xfstt/FAQ
share/doc/xfstt/THANKS
share/doc/xfstt/TODO
share/examples/rc.d/xfstt
share/examples/xfstt/font.properties
@dirrm share/examples/xfstt
@dirrm share/doc/xfstt
