# x11/xfstt/PLIST
@comment $NetBSD: PLIST,v 1.5 2004/01/25 20:46:19 wiz Exp $
bin/xfstt
etc/rc.d/xfstt
man/man1/xfstt.1
share/doc/xfstt/CHANGES
share/doc/xfstt/FAQ
share/doc/xfstt/THANKS
share/doc/xfstt/TODO
share/examples/xfstt/font.properties
@dirrm share/examples/xfstt
@dirrm share/doc/xfstt
