# x11/xfstt/PLIST
@comment $NetBSD: PLIST,v 1.10 2012/06/13 08:23:25 mef Exp $
bin/xfstt
man/man1/xfstt.1
share/doc/xfstt/FAQ
share/doc/xfstt/THANKS
share/doc/xfstt/TODO
share/examples/rc.d/xfstt
share/examples/xfstt/font.properties
@pkgdir lib/X11/fonts/TTF
