# x11/xfstt/PLIST
@comment $NetBSD: PLIST,v 1.4 2004/01/01 20:06:11 salo Exp $
bin/xfstt
etc/rc.d/xfstt
man/man1/xfstt.1
share/doc/xfstt/CHANGES
share/doc/xfstt/FAQ
share/doc/xfstt/THANKS
share/doc/xfstt/TODO
share/examples/xfstt/font.properties
@dirrm share/examples/xfstt
@dirrm share/doc/xfstt
@unexec ${RMDIR} -p %D/lib/X11/fonts/TrueType 2> /dev/null || ${TRUE}
