# x11/xfstt/PLIST
@comment $NetBSD: PLIST,v 1.2 2002/08/22 02:22:14 grant Exp $
bin/xfstt
etc/rc.d/xfstt
man/man1/xfstt.1
share/doc/xfstt/FAQ
share/doc/xfstt/CHANGES
share/doc/xfstt/INSTALL
share/doc/xfstt/THANKS
lib/X11/fonts/TrueType/.keepme
@dirrm lib/X11/fonts/TrueType
@dirrm share/doc/xfstt
