# chat/anope/PLIST
@comment $NetBSD: PLIST,v 1.4 2009/06/14 17:37:42 joerg Exp $
lib/anope/anoperc
lib/anope/languages/cat
lib/anope/languages/de
lib/anope/languages/en_us
lib/anope/languages/es
lib/anope/languages/fr
lib/anope/languages/gr
lib/anope/languages/hun
lib/anope/languages/it
lib/anope/languages/nl
lib/anope/languages/pl
lib/anope/languages/pt
lib/anope/languages/ru
lib/anope/languages/tr
lib/anope/listchans
lib/anope/listnicks
lib/anope/modules/bahamut.so
lib/anope/modules/bs_act.so
lib/anope/modules/bs_assign.so
lib/anope/modules/bs_badwords.so
lib/anope/modules/bs_bot.so
lib/anope/modules/bs_botlist.so
lib/anope/modules/bs_fantasy.so
lib/anope/modules/bs_fantasy_kick.so
lib/anope/modules/bs_fantasy_kickban.so
lib/anope/modules/bs_fantasy_owner.so
lib/anope/modules/bs_fantasy_seen.so
lib/anope/modules/bs_fantasy_unban.so
lib/anope/modules/bs_help.so
lib/anope/modules/bs_info.so
lib/anope/modules/bs_kick.so
lib/anope/modules/bs_say.so
lib/anope/modules/bs_set.so
lib/anope/modules/bs_unassign.so
lib/anope/modules/charybdis.so
lib/anope/modules/cs_access.so
lib/anope/modules/cs_akick.so
lib/anope/modules/cs_appendtopic.so
lib/anope/modules/cs_ban.so
lib/anope/modules/cs_clear.so
lib/anope/modules/cs_drop.so
lib/anope/modules/cs_enforce.so
lib/anope/modules/cs_forbid.so
lib/anope/modules/cs_getkey.so
lib/anope/modules/cs_getpass.so
lib/anope/modules/cs_help.so
lib/anope/modules/cs_identify.so
lib/anope/modules/cs_info.so
lib/anope/modules/cs_invite.so
lib/anope/modules/cs_kick.so
lib/anope/modules/cs_list.so
lib/anope/modules/cs_logout.so
lib/anope/modules/cs_modes.so
lib/anope/modules/cs_register.so
lib/anope/modules/cs_sendpass.so
lib/anope/modules/cs_set.so
lib/anope/modules/cs_status.so
lib/anope/modules/cs_suspend.so
lib/anope/modules/cs_tban.so
lib/anope/modules/cs_topic.so
lib/anope/modules/cs_xop.so
lib/anope/modules/dreamforge.so
lib/anope/modules/enc_md5.so
lib/anope/modules/enc_none.so
lib/anope/modules/enc_old.so
lib/anope/modules/enc_sha1.so
lib/anope/modules/he_help.so
lib/anope/modules/hs_del.so
lib/anope/modules/hs_delall.so
lib/anope/modules/hs_group.so
lib/anope/modules/hs_help.so
lib/anope/modules/hs_list.so
lib/anope/modules/hs_off.so
lib/anope/modules/hs_on.so
lib/anope/modules/hs_request.so
lib/anope/modules/hs_set.so
lib/anope/modules/hs_setall.so
lib/anope/modules/hybrid.so
lib/anope/modules/inspircd10.so
lib/anope/modules/inspircd11.so
lib/anope/modules/ms_cancel.so
lib/anope/modules/ms_check.so
lib/anope/modules/ms_del.so
lib/anope/modules/ms_help.so
lib/anope/modules/ms_info.so
lib/anope/modules/ms_list.so
lib/anope/modules/ms_read.so
lib/anope/modules/ms_rsend.so
lib/anope/modules/ms_send.so
lib/anope/modules/ms_sendall.so
lib/anope/modules/ms_set.so
lib/anope/modules/ms_staff.so
lib/anope/modules/ns_access.so
lib/anope/modules/ns_alist.so
lib/anope/modules/ns_drop.so
lib/anope/modules/ns_forbid.so
lib/anope/modules/ns_getemail.so
lib/anope/modules/ns_getpass.so
lib/anope/modules/ns_ghost.so
lib/anope/modules/ns_group.so
lib/anope/modules/ns_help.so
lib/anope/modules/ns_identify.so
lib/anope/modules/ns_info.so
lib/anope/modules/ns_list.so
lib/anope/modules/ns_logout.so
lib/anope/modules/ns_maxemail.so
lib/anope/modules/ns_recover.so
lib/anope/modules/ns_register.so
lib/anope/modules/ns_release.so
lib/anope/modules/ns_saset.so
lib/anope/modules/ns_sendpass.so
lib/anope/modules/ns_set.so
lib/anope/modules/ns_status.so
lib/anope/modules/ns_suspend.so
lib/anope/modules/ns_update.so
lib/anope/modules/os_admin.so
lib/anope/modules/os_akill.so
lib/anope/modules/os_chankill.so
lib/anope/modules/os_chanlist.so
lib/anope/modules/os_clearmodes.so
lib/anope/modules/os_defcon.so
lib/anope/modules/os_global.so
lib/anope/modules/os_help.so
lib/anope/modules/os_ignore.so
lib/anope/modules/os_info.so
lib/anope/modules/os_jupe.so
lib/anope/modules/os_kick.so
lib/anope/modules/os_logonnews.so
lib/anope/modules/os_mode.so
lib/anope/modules/os_modinfo.so
lib/anope/modules/os_modlist.so
lib/anope/modules/os_modload.so
lib/anope/modules/os_modunload.so
lib/anope/modules/os_noop.so
lib/anope/modules/ns_noop_convert.so
lib/anope/modules/os_oline.so
lib/anope/modules/os_oper.so
lib/anope/modules/os_opernews.so
lib/anope/modules/os_quit.so
lib/anope/modules/os_randomnews.so
lib/anope/modules/os_raw.so
lib/anope/modules/os_reload.so
lib/anope/modules/os_restart.so
lib/anope/modules/os_session.so
lib/anope/modules/os_set.so
lib/anope/modules/os_sgline.so
lib/anope/modules/os_shutdown.so
lib/anope/modules/os_sqline.so
lib/anope/modules/os_staff.so
lib/anope/modules/os_stats.so
lib/anope/modules/os_svsnick.so
lib/anope/modules/os_szline.so
lib/anope/modules/os_umode.so
lib/anope/modules/os_update.so
lib/anope/modules/os_userlist.so
lib/anope/modules/plexus2.so
lib/anope/modules/plexus3.so
lib/anope/modules/ptlink.so
lib/anope/modules/rageircd.so
lib/anope/modules/ratbox.so
lib/anope/modules/shadowircd.so
lib/anope/modules/solidircd.so
lib/anope/modules/ultimate2.so
lib/anope/modules/ultimate3.so
lib/anope/modules/unreal31.so
lib/anope/modules/unreal32.so
lib/anope/modules/viagra.so
lib/anope/mydbgen
lib/anope/services
lib/anope/tools/anopesmtp
lib/anope/tools/db-merger
lib/anope/tools/epona2anope
share/doc/anope/BUGS
share/doc/anope/COPYING
share/doc/anope/DEFCON
share/doc/anope/EVENTS
share/doc/anope/FAQ
share/doc/anope/INSTALL
share/doc/anope/IRCD
share/doc/anope/MODULES
share/doc/anope/MYSQL
share/doc/anope/NEWS
share/doc/anope/OLDCHANGES
share/doc/anope/PROXY
share/doc/anope/README
share/doc/anope/WIN32.txt
share/examples/anope/example.chk
share/examples/anope/example.conf
share/examples/anope/tables.sql
share/examples/rc.d/anope
@unexec ${RMDIR} %D/lib/anope/modules/runtime 2>/dev/null || ${TRUE}
@unexec ${RMDIR} %D/lib/anope/modules 2>/dev/null || ${TRUE}
@unexec ${RMDIR} %D/lib/anope/logs 2>/dev/null || ${TRUE}
@unexec ${RMDIR} %D/lib/anope/backups 2>/dev/null || ${TRUE}
@unexec ${RMDIR} %D/lib/anope 2>/dev/null || ${TRUE}
