# emulators/vmware-module/PLIST
@comment $NetBSD: PLIST,v 1.2 2002/02/15 10:12:37 skrll Exp $
etc/rc.conf.d/linuxrtc
etc/rc.conf.d/vmmon
etc/rc.conf.d/vmnet
etc/rc.conf.d/vmnetif
etc/rc.conf.d/vmware
etc/rc.d/linuxrtc
etc/rc.d/vmmon
etc/rc.d/vmnet
etc/rc.d/vmnetif
lkm/if_hub.o
lkm/if_hub_post.sh
lkm/linuxrtc.o
lkm/linuxrtc_post.sh
lkm/vmmon.o
lkm/vmmon_post.sh
@dirrm etc/rc.conf.d
@unexec ${RM} -f %D/lkm/if_hub. 2>/dev/null || true
@unexec ${RM} -f %D/lkm/linuxrtc 2>/dev/null || true
@unexec ${RM} -f %D/lkm/vmmon 2>/dev/null || true
@unexec ${RMDIR} %D/lkm 2>/dev/null || true
