# mail/dbmail/PLIST
@comment $NetBSD: PLIST,v 1.4 2005/05/03 15:36:38 wiz Exp $
lib/dbmail/lib${SQLDB}dbmail.la
lib/dbmail/libauthdbmail.la
lib/dbmail/libsortdbmail.la
lib/libdbmail.la
man/man1/dbmail-smtp.1
man/man8/dbmail-imapd.8
man/man8/dbmail-lmtpd.8
man/man8/dbmail-pop3d.8
man/man8/dbmail-users.8
man/man8/dbmail-util.8
sbin/dbmail-imapd
sbin/dbmail-lmtpd
sbin/dbmail-pop3d
sbin/dbmail-smtp
sbin/dbmail-users
sbin/dbmail-util
share/examples/dbmail/dbmail.conf
share/examples/rc.d/dbmail
share/examples/rc.d/dbmailimapd
share/examples/rc.d/dbmaillmtpd
share/examples/rc.d/dbmailpop3d
@dirrm share/examples/dbmail
@dirrm lib/dbmail
