# security/mit-krb5/PLIST
@comment $NetBSD: PLIST,v 1.10 2005/05/03 15:49:11 wiz Exp $
bin/compile_et
bin/gss-client
bin/kdestroy
bin/kftp
bin/kinit
bin/klist
bin/kpasswd
bin/krb5-config
bin/${KRB5_PREFIX}rcp
bin/${KRB5_PREFIX}rlogin
bin/${KRB5_PREFIX}rsh
bin/ksu
bin/${KRB5_PREFIX}telnet
bin/kvno
bin/sclient
bin/sim_client
bin/uuclient
include/com_err.h
include/gssapi/gssapi.h
include/gssapi/gssapi_generic.h
include/gssapi/gssapi_krb5.h
include/gssrpc/auth.h
include/gssrpc/auth_gss.h
include/gssrpc/auth_gssapi.h
include/gssrpc/auth_unix.h
include/gssrpc/clnt.h
include/gssrpc/netdb.h
include/gssrpc/pmap_clnt.h
include/gssrpc/pmap_prot.h
include/gssrpc/pmap_rmt.h
include/gssrpc/rename.h
include/gssrpc/rpc.h
include/gssrpc/rpc_msg.h
include/gssrpc/svc.h
include/gssrpc/svc_auth.h
include/gssrpc/types.h
include/gssrpc/xdr.h
include/krb5.h
include/profile.h
lib/libcom_err.la
lib/libdes425.la
lib/libgssapi_krb5.la
lib/libgssrpc.la
lib/libk5crypto.la
lib/libkadm5clnt.la
lib/libkadm5srv.la
lib/libkdb5.la
lib/libkrb5.la
lib/libkrb5support.la
man/man1/kdestroy.1
man/man1/kerberos.1
man/man1/kftp.1
man/man1/kinit.1
man/man1/klist.1
man/man1/kpasswd.1
man/man1/krb5-config.1
man/man1/krb5-send-pr.1
man/man1/${KRB5_PREFIX}rcp.1
man/man1/${KRB5_PREFIX}rlogin.1
man/man1/${KRB5_PREFIX}rsh.1
man/man1/ksu.1
man/man1/${KRB5_PREFIX}telnet.1
man/man1/kvno.1
man/man1/sclient.1
man/man5/.k5login.5
man/man5/kdc.conf.5
man/man5/krb5.conf.5
man/man8/k5srvutil.8
man/man8/kadmin.8
man/man8/kadmin.local.8
man/man8/kadmind.8
man/man8/kdb5_util.8
man/man8/kftpd.8
man/man8/klogind.8
man/man8/kprop.8
man/man8/kpropd.8
man/man8/krb5kdc.8
man/man8/kshd.8
man/man8/${KRB5_PREFIX}telnetd.8
man/man8/ktutil.8
man/man8/login.krb5.8
man/man8/sserver.8
sbin/gss-server
sbin/k5srvutil
sbin/kadmin
sbin/kadmin.local
sbin/kadmind
sbin/kdb5_util
sbin/kftpd
sbin/klogind
sbin/kprop
sbin/kpropd
sbin/krb5-send-pr
sbin/krb5kdc
sbin/kshd
sbin/${KRB5_PREFIX}telnetd
sbin/ktutil
sbin/login.krb5
sbin/sim_server
sbin/sserver
sbin/uuserver
share/et/et_c.awk
share/et/et_h.awk
share/examples/krb5/kdc.conf
share/examples/krb5/krb5.conf
share/examples/krb5/services.append
share/examples/rc.d/kadmind
share/examples/rc.d/kdc
share/gnats/mit
@dirrm share/gnats
@dirrm share/examples/krb5
@dirrm share/et
@dirrm include/gssrpc
@dirrm include/gssapi
