# security/stunnel/PLIST
@comment $NetBSD: PLIST,v 1.7 2005/05/02 20:34:05 reed Exp $
lib/libstunnel.la
man/man8/stunnel.8
man/man8/stunnel.fr.8
man/man8/stunnel.pl.8
sbin/stunnel
sbin/stunnel3
share/doc/stunnel/AUTHORS
share/doc/stunnel/BUGS
share/doc/stunnel/COPYING
share/doc/stunnel/COPYRIGHT.GPL
share/doc/stunnel/CREDITS
share/doc/stunnel/ChangeLog
share/doc/stunnel/INSTALL
share/doc/stunnel/INSTALL.W32
share/doc/stunnel/NEWS
share/doc/stunnel/PORTS
share/doc/stunnel/README
share/doc/stunnel/TODO
share/doc/stunnel/examples/ca.html
share/doc/stunnel/examples/ca.pl
share/doc/stunnel/examples/importCA.html
share/doc/stunnel/examples/importCA.sh
share/doc/stunnel/examples/script.sh
share/doc/stunnel/examples/stunnel.init
share/doc/stunnel/examples/stunnel.spec
share/doc/stunnel/stunnel.fr.html
share/doc/stunnel/stunnel.html
share/doc/stunnel/stunnel.pl.html
share/examples/rc.d/stunnel
share/examples/stunnel/stunnel.conf-sample
@dirrm share/examples/stunnel
@dirrm share/doc/stunnel/examples
@dirrm share/doc/stunnel
