# chat/bitlbee/PLIST
@comment $NetBSD: PLIST,v 1.6 2009/07/15 22:14:18 tonio Exp $
man/man5/bitlbee.conf.5
man/man8/bitlbee.8
sbin/bitlbee
share/bitlbee/help.txt
share/bitlbee/motd.txt
share/examples/rc.d/bitlbee
share/doc/bitlbee/AUTHORS
share/doc/bitlbee/CREDITS
share/doc/bitlbee/FAQ
share/doc/bitlbee/README
share/doc/bitlbee/user-guide.html
share/doc/bitlbee/user-guide.txt
share/examples/bitlbee/bitlbee.conf
