# security/openssh/files/smf/manifest.xml
<?xml version='1.0'?>
<!DOCTYPE service_bundle SYSTEM '/usr/share/lib/xml/dtd/service_bundle.dtd.1'>
<service_bundle type='manifest' name='export'>
  <service name='@SMF_PREFIX@/@SMF_NAME@' type='service' version='1'>
    <create_default_instance enabled='false'/>
    <single_instance/>
    <dependency name='fs-local' grouping='require_all' restart_on='none' type='service'>
      <service_fmri value='svc:/system/filesystem/local'/>
    </dependency>
    <dependency name='net-loopback' grouping='require_all' restart_on='none' type='service'>
      <service_fmri value='svc:/network/loopback'/>
    </dependency>
    <dependency name='net-physical' grouping='require_all' restart_on='none' type='service'>
      <service_fmri value='svc:/network/physical'/>
    </dependency>
    <dependency name='cryptosvc' grouping='require_all' restart_on='none' type='service'>
      <service_fmri value='svc:/system/cryptosvc'/>
    </dependency>
    <dependency name='utmp' grouping='require_all' restart_on='none' type='service'>
      <service_fmri value='svc:/system/utmp'/>
    </dependency>
    <dependency name='config_data' grouping='require_all' restart_on='restart' type='path'>
      <service_fmri value='file://localhost@PKG_SYSCONFDIR@/sshd_config'/>
    </dependency>
    <dependent name='openssh_multi-user-server' restart_on='none' grouping='optional_all'>
      <service_fmri value='svc:/milestone/multi-user-server'/>
    </dependent>
    <exec_method name='start' type='method' exec='@PREFIX@/@SMF_METHOD_FILE.sshd@ start' timeout_seconds='60'/>
    <exec_method name='stop' type='method' exec=':kill' timeout_seconds='60'/>
    <exec_method name='refresh' type='method' exec='@PREFIX@/@SMF_METHOD_FILE.sshd@ restart' timeout_seconds='60'/>
    <property_group name='general' type='framework'>
      <property name='action_authorization' type='astring'/>
    </property_group>
    <property_group name='startd' type='framework'>
      <propval name='ignore_error' type='astring' value='core,signal'/>
    </property_group>
    <template>
      <common_name>
        <loctext xml:lang='C'>OpenSSH server</loctext>
      </common_name>
      <documentation>
        <manpage title='sshd' section='1M' manpath='@PREFIX@/@PKGMANDIR@'/>
      </documentation>
    </template>
  </service>
</service_bundle>
