# net/aiccu/PLIST
@comment $NetBSD: PLIST,v 1.3 2009/06/14 18:09:20 joerg Exp $
man/man1/aiccu.1
sbin/aiccu
share/doc/aiccu/COPYRIGHT
share/doc/aiccu/HOWTO
share/doc/aiccu/LICENSE
share/doc/aiccu/README
share/examples/aiccu/aiccu.conf
share/examples/rc.d/aiccu
