# net/aiccu/PLIST
@comment $NetBSD: PLIST,v 1.2 2007/01/18 14:42:17 ghen Exp $
man/man1/aiccu.1
sbin/aiccu
share/doc/aiccu/COPYRIGHT
share/doc/aiccu/HOWTO
share/doc/aiccu/LICENSE
share/doc/aiccu/README
share/examples/aiccu/aiccu.conf
share/examples/rc.d/aiccu
@dirrm share/examples/aiccu
@dirrm share/doc/aiccu
