# comms/conserver/PLIST
@comment $NetBSD: PLIST,v 1.3 2003/05/22 06:35:02 wiz Exp $
bin/console
etc/rc.d/conserver
man/man1/console.1
man/man5/conserver.cf.5
man/man5/conserver.passwd.5
man/man8/conserver.8
sbin/conserver
share/examples/conserver/conserver.cf
share/examples/conserver/conserver.passwd
@dirrm share/examples/conserver
