# pkgtools/rc.subr/PLIST
@comment $NetBSD: PLIST,v 1.1.1.1 2003/06/19 16:07:10 grant Exp $
etc/defaults/rc.conf
etc/rc.conf.example
etc/rc.subr
etc/rc.d/DAEMON
etc/rc.d/LOGIN
etc/rc.d/NETWORKING
etc/rc.d/SERVERS
@unexec ${RMDIR} %Detc/rc.d 2> /dev/null || ${TRUE}
@dirrm etc/defaults
