# sysutils/collectd/PLIST
@comment $NetBSD: PLIST,v 1.1 2014/02/17 11:21:54 fhajny Exp $
bin/collectd-nagios
bin/collectd-tg
bin/collectdctl
include/collectd/client.h
include/collectd/lcc_features.h
include/collectd/network.h
include/collectd/network_buffer.h
lib/collectd/aggregation.la
lib/collectd/apcups.la
${PLIST.apple-sensors}lib/collectd/apple_sensors.la
${PLIST.battery}lib/collectd/battery.la
${PLIST.contextswitch}lib/collectd/contextswitch.la
${PLIST.cpu}lib/collectd/cpu.la
lib/collectd/csv.la
${PLIST.df}lib/collectd/df.la
${PLIST.disk}lib/collectd/disk.la
lib/collectd/email.la
lib/collectd/exec.la
lib/collectd/filecount.la
lib/collectd/hddtemp.la
${PLIST.interface}lib/collectd/interface.la
${PLIST.load}lib/collectd/load.la
lib/collectd/logfile.la
lib/collectd/match_empty_counter.la
lib/collectd/match_hashed.la
lib/collectd/match_regex.la
lib/collectd/match_timediff.la
lib/collectd/match_value.la
lib/collectd/mbmon.la
${PLIST.memory}lib/collectd/memory.la
${PLIST.nfs}lib/collectd/nfs.la
lib/collectd/ntpd.la
lib/collectd/olsrd.la
lib/collectd/openvpn.la
${PLIST.perl}lib/collectd/perl.la
${PLIST.pf}lib/collectd/pf.la
lib/collectd/powerdns.la
${PLIST.processes}lib/collectd/processes.la
${PLIST.python}lib/collectd/python.la
lib/collectd/statsd.la
${PLIST.swap}lib/collectd/swap.la
${PLIST.syslog}lib/collectd/syslog.la
lib/collectd/table.la
lib/collectd/tail.la
lib/collectd/tail_csv.la
lib/collectd/target_notification.la
lib/collectd/target_replace.la
lib/collectd/target_scale.la
lib/collectd/target_set.la
lib/collectd/target_v5upgrade.la
lib/collectd/teamspeak2.la
${PLIST.tcpconns}lib/collectd/tcpconns.la
lib/collectd/threshold.la
lib/collectd/unixsock.la
${PLIST.uptime}lib/collectd/uptime.la
${PLIST.users}lib/collectd/users.la
lib/collectd/uuid.la
lib/collectd/write_graphite.la
${PLIST.zfs-arc}lib/collectd/zfs_arc.la
lib/libcollectdclient.la
lib/pkgconfig/libcollectdclient.pc
man/man1/collectd-nagios.1
man/man1/collectd-tg.1
man/man1/collectd.1
man/man1/collectdctl.1
man/man1/collectdmon.1
man/man5/collectd-email.5
man/man5/collectd-exec.5
man/man5/collectd-java.5
man/man5/collectd-perl.5
man/man5/collectd-python.5
man/man5/collectd-snmp.5
man/man5/collectd-threshold.5
man/man5/collectd-unixsock.5
man/man5/collectd.conf.5
man/man5/types.db.5
sbin/collectd
sbin/collectdmon
share/collectd/postgresql_default.conf
share/collectd/types.db
share/examples/collectd/collectd.conf
share/examples/rc.d/collectd
