# net/Geomyidae/PLIST
@comment $NetBSD: PLIST,v 1.1.1.1 2008/04/06 11:50:36 obache Exp $
man/man8/Geomyidae.8
sbin/Geomyidae
share/doc/Geomyidae/LICENSE
share/doc/Geomyidae/README
share/examples/Geomyidae/index.gph
share/examples/rc.d/Geomyidae
@dirrm share/examples/Geomyidae
@dirrm share/doc/Geomyidae
