# net/Geomyidae/PLIST
@comment $NetBSD: PLIST,v 1.2 2009/06/14 18:09:20 joerg Exp $
man/man8/Geomyidae.8
sbin/Geomyidae
share/doc/Geomyidae/LICENSE
share/doc/Geomyidae/README
share/examples/Geomyidae/index.gph
share/examples/rc.d/Geomyidae
