# net/gnugk/PLIST
@comment $NetBSD: PLIST,v 1.3 2005/05/02 20:34:02 reed Exp $
bin/gnugk
share/doc/gnugk/GkStatusMsgFmt.html
share/doc/gnugk/changing_ip.txt
share/doc/gnugk/comparison.txt
share/doc/gnugk/compiling.txt
share/doc/gnugk/copying
share/doc/gnugk/gk_faq.html
share/doc/gnugk/gkstatus.txt
share/doc/gnugk/manual.sgml
share/doc/gnugk/performance.txt
share/doc/gnugk/reference.txt
share/doc/gnugk/signalling.txt
share/doc/gnugk/test-status.txt
share/doc/gnugk/trace_levels.txt
share/doc/gnugk/tutorial.txt
share/examples/gnugk/billing.pl
share/examples/gnugk/child.ini
share/examples/gnugk/complete.ini
share/examples/gnugk/fileacct.ini
share/examples/gnugk/mysql.ini
share/examples/gnugk/neighbor.ini
share/examples/gnugk/proxy.ini
share/examples/gnugk/radacct.ini
share/examples/gnugk/radaliasauth.ini
share/examples/gnugk/radauth.ini
share/examples/gnugk/routed.ini
share/examples/gnugk/rrqauth.ini
share/examples/gnugk/sql.ini
share/examples/gnugk/status.ini
share/examples/gnugk/unregister_ip.pl
share/examples/gnugk/vqueue.ini
share/examples/rc.d/gnugk
@dirrm share/examples/gnugk
@dirrm share/doc/gnugk
