# inputmethod/ja-freewnn-server-bin/PLIST
@comment $NetBSD: PLIST,v 1.4 2005/05/05 19:37:23 kei Exp $
man/ja_JP.EUC/man1/atod.1
man/ja_JP.EUC/man1/atof.1
man/ja_JP.EUC/man1/dtoa.1
man/ja_JP.EUC/man1/jserver.1
man/ja_JP.EUC/man1/oldatonewa.1
man/ja_JP.EUC/man1/uum.1
man/ja_JP.EUC/man1/wddel.1
man/ja_JP.EUC/man1/wdreg.1
man/ja_JP.EUC/man1/wnnkill.1
man/ja_JP.EUC/man1/wnnstat.1
man/ja_JP.EUC/man1/wnntouch.1
man/ja_JP.EUC/man4/2a_ctrl.4
man/ja_JP.EUC/man4/2b_romkana.4
man/ja_JP.EUC/man4/cvt_key_tbl.4
man/ja_JP.EUC/man4/fzk.data.4
man/ja_JP.EUC/man4/fzk.u.4
man/ja_JP.EUC/man4/hinsi_data.4
man/ja_JP.EUC/man4/jserverrc.4
man/ja_JP.EUC/man4/mode.4
man/ja_JP.EUC/man4/serverdefs.4
man/ja_JP.EUC/man4/ujis_dic.4
man/ja_JP.EUC/man4/uumkey.4
man/ja_JP.EUC/man4/uumrc.4
man/ja_JP.EUC/man4/wnnenvrc.4
man/ja_JP.EUC/man5/pubdic.5
man/ja_JP.EUC/man5/usr_dic.5
sbin/atod
sbin/atof
sbin/dtoa
sbin/jserver
sbin/oldatonewa
sbin/wddel
sbin/wdreg
sbin/wnnkill
sbin/wnnstat
sbin/wnntouch
share/examples/rc.d/ja_freewnn
share/wnn/ja_JP/jserverrc
