# mail/amavisd-new/files/smf/manifest.xml
<?xml version="1.0"?>
<!DOCTYPE service_bundle SYSTEM "/usr/share/lib/xml/dtd/service_bundle.dtd.1">
<service_bundle type="manifest" name="export">
  <service name="@SMF_PREFIX@/@SMF_NAME@" type="service" version="1">
    <create_default_instance enabled="false" />
    <single_instance />
    <dependency name='fs-local' grouping='require_all' restart_on='none' type='service'>
      <service_fmri value='svc:/system/filesystem/local' />
    </dependency>
    <dependency name='loopback' grouping='require_all' restart_on='error' type='service'>
      <service_fmri value='svc:/network/loopback:default' />
    </dependency>
    <dependency name='physical' grouping='require_all' restart_on='error' type='service'>
      <service_fmri value='svc:/network/physical:default' />
    </dependency>
    <dependency name='name-services' grouping='require_all'  restart_on='refresh' type='service'>
      <service_fmri value='svc:/milestone/name-services' />
    </dependency>
    <dependency name='system-log' grouping='optional_all' restart_on='none' type='service'>
      <service_fmri value='svc:/system/system-log' />
    </dependency>
    <exec_method type="method" name="start" exec="@PREFIX@/sbin/amavisd start" timeout_seconds="60" />
    <exec_method type="method" name="stop" exec="@PREFIX@/sbin/amavisd stop" timeout_seconds="120" />
    <exec_method type="method" name="refresh" exec="@PREFIX@/sbin/amavisd reload" timeout_seconds="60" />
    <property_group name="startd" type="framework">
      <propval name="duration" type="astring" value="contract" />
      <propval name="ignore_error" type="astring" value="core,signal" />
    </property_group>
    <property_group name="application" type="application">
      <propval name="config_file" type="astring" value="@PKG_SYSCONFDIR@/amavisd.conf" />
    </property_group>
    <template>
      <common_name>
        <loctext xml:lang="C">Amavisd</loctext>
      </common_name>
    </template>
  </service>
</service_bundle>
