# net/dhid/pkg/PLIST
@comment $NetBSD: PLIST,v 1.2 1999/08/04 14:04:37 jlam Exp $
etc/rc.d/dhid.sh
libexec/dhid
man/man1/dtool.1
man/man8/dhid.8
sbin/dtool
share/doc/dhid/COPYRIGHT
share/doc/dhid/README
share/examples/dhid.conf
@dirrm share/doc/dhid
