# mail/tmda/PLIST
@comment $NetBSD: PLIST,v 1.12 2006/11/07 07:42:09 schmonz Exp $
bin/tmda-address
bin/tmda-check-address
bin/tmda-filter
bin/tmda-inject
bin/tmda-keygen
bin/tmda-ofmipd
bin/tmda-pending
bin/tmda-rfilter
bin/tmda-sendmail
${PYSITELIB}/TMDA/Address.py
${PYSITELIB}/TMDA/Address.pyc
${PYSITELIB}/TMDA/Auth.py
${PYSITELIB}/TMDA/Auth.pyc
${PYSITELIB}/TMDA/AutoResponse.py
${PYSITELIB}/TMDA/AutoResponse.pyc
${PYSITELIB}/TMDA/Cookie.py
${PYSITELIB}/TMDA/Cookie.pyc
${PYSITELIB}/TMDA/Defaults.py
${PYSITELIB}/TMDA/Defaults.pyc
${PYSITELIB}/TMDA/Deliver.py
${PYSITELIB}/TMDA/Deliver.pyc
${PYSITELIB}/TMDA/Errors.py
${PYSITELIB}/TMDA/Errors.pyc
${PYSITELIB}/TMDA/FilterParser.py
${PYSITELIB}/TMDA/FilterParser.pyc
${PYSITELIB}/TMDA/HMAC.py
${PYSITELIB}/TMDA/HMAC.pyc
${PYSITELIB}/TMDA/MTA.py
${PYSITELIB}/TMDA/MTA.pyc
${PYSITELIB}/TMDA/MessageLogger.py
${PYSITELIB}/TMDA/MessageLogger.pyc
${PYSITELIB}/TMDA/Pending.py
${PYSITELIB}/TMDA/Pending.pyc
${PYSITELIB}/TMDA/SMTP.py
${PYSITELIB}/TMDA/SMTP.pyc
${PYSITELIB}/TMDA/Util.py
${PYSITELIB}/TMDA/Util.pyc
${PYSITELIB}/TMDA/Version.py
${PYSITELIB}/TMDA/Version.pyc
${PYSITELIB}/TMDA/__init__.py
${PYSITELIB}/TMDA/__init__.pyc
${PYSITELIB}/TMDA/pythonlib/email/__init__.py
${PYSITELIB}/TMDA/pythonlib/email/__init__.pyc
${PYSITELIB}/TMDA/pythonlib/email/_parseaddr.py
${PYSITELIB}/TMDA/pythonlib/email/_parseaddr.pyc
${PYSITELIB}/TMDA/pythonlib/email/base64mime.py
${PYSITELIB}/TMDA/pythonlib/email/base64mime.pyc
${PYSITELIB}/TMDA/pythonlib/email/charset.py
${PYSITELIB}/TMDA/pythonlib/email/charset.pyc
${PYSITELIB}/TMDA/pythonlib/email/encoders.py
${PYSITELIB}/TMDA/pythonlib/email/encoders.pyc
${PYSITELIB}/TMDA/pythonlib/email/errors.py
${PYSITELIB}/TMDA/pythonlib/email/errors.pyc
${PYSITELIB}/TMDA/pythonlib/email/feedparser.py
${PYSITELIB}/TMDA/pythonlib/email/feedparser.pyc
${PYSITELIB}/TMDA/pythonlib/email/generator.py
${PYSITELIB}/TMDA/pythonlib/email/generator.pyc
${PYSITELIB}/TMDA/pythonlib/email/header.py
${PYSITELIB}/TMDA/pythonlib/email/header.pyc
${PYSITELIB}/TMDA/pythonlib/email/iterators.py
${PYSITELIB}/TMDA/pythonlib/email/iterators.pyc
${PYSITELIB}/TMDA/pythonlib/email/message.py
${PYSITELIB}/TMDA/pythonlib/email/message.pyc
${PYSITELIB}/TMDA/pythonlib/email/mime/__init__.py
${PYSITELIB}/TMDA/pythonlib/email/mime/__init__.pyc
${PYSITELIB}/TMDA/pythonlib/email/mime/application.py
${PYSITELIB}/TMDA/pythonlib/email/mime/application.pyc
${PYSITELIB}/TMDA/pythonlib/email/mime/audio.py
${PYSITELIB}/TMDA/pythonlib/email/mime/audio.pyc
${PYSITELIB}/TMDA/pythonlib/email/mime/base.py
${PYSITELIB}/TMDA/pythonlib/email/mime/base.pyc
${PYSITELIB}/TMDA/pythonlib/email/mime/image.py
${PYSITELIB}/TMDA/pythonlib/email/mime/image.pyc
${PYSITELIB}/TMDA/pythonlib/email/mime/message.py
${PYSITELIB}/TMDA/pythonlib/email/mime/message.pyc
${PYSITELIB}/TMDA/pythonlib/email/mime/multipart.py
${PYSITELIB}/TMDA/pythonlib/email/mime/multipart.pyc
${PYSITELIB}/TMDA/pythonlib/email/mime/nonmultipart.py
${PYSITELIB}/TMDA/pythonlib/email/mime/nonmultipart.pyc
${PYSITELIB}/TMDA/pythonlib/email/mime/text.py
${PYSITELIB}/TMDA/pythonlib/email/mime/text.pyc
${PYSITELIB}/TMDA/pythonlib/email/parser.py
${PYSITELIB}/TMDA/pythonlib/email/parser.pyc
${PYSITELIB}/TMDA/pythonlib/email/quoprimime.py
${PYSITELIB}/TMDA/pythonlib/email/quoprimime.pyc
${PYSITELIB}/TMDA/pythonlib/email/utils.py
${PYSITELIB}/TMDA/pythonlib/email/utils.pyc
${HTMLDIR}/config-vars.html
${HTMLDIR}/img/blocksmall.png
${HTMLDIR}/img/stripes.gif
${HTMLDIR}/img/tmda-186x253.png
${HTMLDIR}/img/tmda-47x64.png
${HTMLDIR}/img/tmda-562x765.png
${HTMLDIR}/img/tmda-70x95.png
${HTMLDIR}/img/tmda-93x127.png
${HTMLDIR}/img/tmda-nt-616x616.png
${HTMLDIR}/index.html
${HTMLDIR}/tmda.css
${DOCDIR}/CODENAMES
${DOCDIR}/COPYING
${DOCDIR}/CRYPTO
${DOCDIR}/ChangeLog
${DOCDIR}/INSTALL
${DOCDIR}/NEWS
${DOCDIR}/README
${DOCDIR}/THANKS
${DOCDIR}/UPGRADE
${EGDIR}/ChangeLog
${EGDIR}/collectaddys
${EGDIR}/def2html
${EGDIR}/printcdb
${EGDIR}/printdbm
${EGDIR}/sendit.sh
${EGDIR}/smtp-check-sender
${EGDIR}/tmda.el
${EGDIR}/tmda.spec
${EGDIR}/tofmipd.init
${EGDIR}/tofmipd.sysconfig
${EGDIR}/update-internaldomains
${EGDIR}/vadduser-tmda
${EGDIR}/vmailmgr-vdir.sh
${EGDIR}/vpopmail-vdir.sh
${EGDIR}/vtmdarc
${EGDIR}/wrapfd3.sh
${EGDIR}/dot-tmda/README
${EGDIR}/dot-tmda/config
${EGDIR}/dot-tmda/crypt_key
${EGDIR}/dot-tmda/filters/incoming
${EGDIR}/dot-tmda/filters/outgoing
${EGDIR}/dot-tmda/lists/blacklist
${EGDIR}/dot-tmda/lists/confirmed
${EGDIR}/dot-tmda/lists/whitelist
${EGDIR}/ofmipd-stunnel-xinetd/README.txt
${EGDIR}/ofmipd-stunnel-xinetd/stunnel-wrapper
${EGDIR}/ofmipd-stunnel-xinetd/stunnel.conf
${EGDIR}/ofmipd-stunnel-xinetd/tmda-ofmipd-wrapper
${EGDIR}/ofmipd-stunnel-xinetd/xinetd.d-tmda-ofmipd-starttls
share/examples/rc.d/tofmipd
share/tmda/bounce_expired_dated.txt
share/tmda/bounce_fail_dated.txt
share/tmda/bounce_fail_keyword.txt
share/tmda/bounce_fail_sender.txt
share/tmda/bounce_incoming.txt
share/tmda/bounce_invalid_confirmation.txt
share/tmda/bounce_missing_pending.txt
share/tmda/confirm_accept.txt
share/tmda/confirm_request.txt
@dirrm share/tmda
@dirrm ${DOCDIR}
@dirrm ${HTMLDIR}/img
@dirrm ${HTMLDIR}
@dirrm ${EGDIR}/ofmipd-stunnel-xinetd
@dirrm ${EGDIR}/dot-tmda/lists
@dirrm ${EGDIR}/dot-tmda/filters
@dirrm ${EGDIR}/dot-tmda
@dirrm ${EGDIR}
@dirrm ${PYSITELIB}/TMDA/pythonlib/email/mime
@dirrm ${PYSITELIB}/TMDA/pythonlib/email
@dirrm ${PYSITELIB}/TMDA/pythonlib
@dirrm ${PYSITELIB}/TMDA
