# mail/tmda/PLIST
@comment $NetBSD: PLIST,v 1.11 2005/05/02 20:34:01 reed Exp $
bin/tmda-address
bin/tmda-check-address
bin/tmda-filter
bin/tmda-inject
bin/tmda-keygen
bin/tmda-ofmipd
bin/tmda-pending
bin/tmda-rfilter
bin/tmda-sendmail
${PYSITELIB}/TMDA/Address.py
${PYSITELIB}/TMDA/Address.pyc
${PYSITELIB}/TMDA/Auth.py
${PYSITELIB}/TMDA/Auth.pyc
${PYSITELIB}/TMDA/AutoResponse.py
${PYSITELIB}/TMDA/AutoResponse.pyc
${PYSITELIB}/TMDA/Cookie.py
${PYSITELIB}/TMDA/Cookie.pyc
${PYSITELIB}/TMDA/Defaults.py
${PYSITELIB}/TMDA/Defaults.pyc
${PYSITELIB}/TMDA/Deliver.py
${PYSITELIB}/TMDA/Deliver.pyc
${PYSITELIB}/TMDA/Errors.py
${PYSITELIB}/TMDA/Errors.pyc
${PYSITELIB}/TMDA/FilterParser.py
${PYSITELIB}/TMDA/FilterParser.pyc
${PYSITELIB}/TMDA/HMAC.py
${PYSITELIB}/TMDA/HMAC.pyc
${PYSITELIB}/TMDA/MTA.py
${PYSITELIB}/TMDA/MTA.pyc
${PYSITELIB}/TMDA/MessageLogger.py
${PYSITELIB}/TMDA/MessageLogger.pyc
${PYSITELIB}/TMDA/Pending.py
${PYSITELIB}/TMDA/Pending.pyc
${PYSITELIB}/TMDA/SMTP.py
${PYSITELIB}/TMDA/SMTP.pyc
${PYSITELIB}/TMDA/Util.py
${PYSITELIB}/TMDA/Util.pyc
${PYSITELIB}/TMDA/Version.py
${PYSITELIB}/TMDA/Version.pyc
${PYSITELIB}/TMDA/__init__.py
${PYSITELIB}/TMDA/__init__.pyc
${PYSITELIB}/TMDA/pythonlib/email/Charset.py
${PYSITELIB}/TMDA/pythonlib/email/Charset.pyc
${PYSITELIB}/TMDA/pythonlib/email/Encoders.py
${PYSITELIB}/TMDA/pythonlib/email/Encoders.pyc
${PYSITELIB}/TMDA/pythonlib/email/Errors.py
${PYSITELIB}/TMDA/pythonlib/email/Errors.pyc
${PYSITELIB}/TMDA/pythonlib/email/Generator.py
${PYSITELIB}/TMDA/pythonlib/email/Generator.pyc
${PYSITELIB}/TMDA/pythonlib/email/Header.py
${PYSITELIB}/TMDA/pythonlib/email/Header.pyc
${PYSITELIB}/TMDA/pythonlib/email/Iterators.py
${PYSITELIB}/TMDA/pythonlib/email/Iterators.pyc
${PYSITELIB}/TMDA/pythonlib/email/MIMEAudio.py
${PYSITELIB}/TMDA/pythonlib/email/MIMEAudio.pyc
${PYSITELIB}/TMDA/pythonlib/email/MIMEBase.py
${PYSITELIB}/TMDA/pythonlib/email/MIMEBase.pyc
${PYSITELIB}/TMDA/pythonlib/email/MIMEImage.py
${PYSITELIB}/TMDA/pythonlib/email/MIMEImage.pyc
${PYSITELIB}/TMDA/pythonlib/email/MIMEMessage.py
${PYSITELIB}/TMDA/pythonlib/email/MIMEMessage.pyc
${PYSITELIB}/TMDA/pythonlib/email/MIMEMultipart.py
${PYSITELIB}/TMDA/pythonlib/email/MIMEMultipart.pyc
${PYSITELIB}/TMDA/pythonlib/email/MIMENonMultipart.py
${PYSITELIB}/TMDA/pythonlib/email/MIMENonMultipart.pyc
${PYSITELIB}/TMDA/pythonlib/email/MIMEText.py
${PYSITELIB}/TMDA/pythonlib/email/MIMEText.pyc
${PYSITELIB}/TMDA/pythonlib/email/Message.py
${PYSITELIB}/TMDA/pythonlib/email/Message.pyc
${PYSITELIB}/TMDA/pythonlib/email/Parser.py
${PYSITELIB}/TMDA/pythonlib/email/Parser.pyc
${PYSITELIB}/TMDA/pythonlib/email/Utils.py
${PYSITELIB}/TMDA/pythonlib/email/Utils.pyc
${PYSITELIB}/TMDA/pythonlib/email/__init__.py
${PYSITELIB}/TMDA/pythonlib/email/__init__.pyc
${PYSITELIB}/TMDA/pythonlib/email/_compat21.py
${PYSITELIB}/TMDA/pythonlib/email/_compat21.pyc
${PYSITELIB}/TMDA/pythonlib/email/_compat22.py
${PYSITELIB}/TMDA/pythonlib/email/_compat22.pyc
${PYSITELIB}/TMDA/pythonlib/email/_parseaddr.py
${PYSITELIB}/TMDA/pythonlib/email/_parseaddr.pyc
${PYSITELIB}/TMDA/pythonlib/email/base64MIME.py
${PYSITELIB}/TMDA/pythonlib/email/base64MIME.pyc
${PYSITELIB}/TMDA/pythonlib/email/quopriMIME.py
${PYSITELIB}/TMDA/pythonlib/email/quopriMIME.pyc
${HTMLDIR}/challengeresponse.html
${HTMLDIR}/config-client.html
${HTMLDIR}/config-filter.html
${HTMLDIR}/config-pre.html
${HTMLDIR}/config-server.html
${HTMLDIR}/config-vars.html
${HTMLDIR}/config.html
${HTMLDIR}/donations.html
${HTMLDIR}/download.html
${HTMLDIR}/features.html
${HTMLDIR}/filter-sources.html
${HTMLDIR}/history.html
${HTMLDIR}/howto-template.html
${HTMLDIR}/img/tmda-186x253.png
${HTMLDIR}/img/tmda-47x64.png
${HTMLDIR}/img/tmda-562x765.png
${HTMLDIR}/img/tmda-70x95.png
${HTMLDIR}/img/tmda-93x127.png
${HTMLDIR}/img/tmda-nt-616x616.png
${HTMLDIR}/index.html
${HTMLDIR}/install.html
${HTMLDIR}/logos.html
${HTMLDIR}/requirements.html
${HTMLDIR}/tmda-ofmipd.html
${HTMLDIR}/tmda-vdomains.html
${HTMLDIR}/trouble.html
${HTMLDIR}/upgrade.html
${DOCDIR}/CODENAMES
${DOCDIR}/COPYING
${DOCDIR}/CRYPTO
${DOCDIR}/ChangeLog
${DOCDIR}/INSTALL
${DOCDIR}/README
${DOCDIR}/THANKS
${DOCDIR}/UPGRADE
${DOCDIR}/contrib/ChangeLog
${DOCDIR}/contrib/collectaddys
${DOCDIR}/contrib/def2html
${DOCDIR}/contrib/printcdb
${DOCDIR}/contrib/printdbm
${DOCDIR}/contrib/sample.config
${DOCDIR}/contrib/sendit.sh
${DOCDIR}/contrib/smtp-check-sender
${DOCDIR}/contrib/tmda.el
${DOCDIR}/contrib/tmda.spec
${DOCDIR}/contrib/tofmipd.init
${DOCDIR}/contrib/tofmipd.sysconfig
${DOCDIR}/contrib/update-internaldomains
${DOCDIR}/contrib/vadduser-tmda
${DOCDIR}/contrib/vmailmgr-vdir.sh
${DOCDIR}/contrib/vpopmail-vdir.sh
${DOCDIR}/contrib/vtmdarc
${DOCDIR}/contrib/wrapfd3.sh
${SHAREDIR}/bounce.txt
${SHAREDIR}/confirm_accept.txt
${SHAREDIR}/confirm_request.txt
share/examples/rc.d/tofmipd
@dirrm ${SHAREDIR}
@dirrm ${DOCDIR}/contrib
@dirrm ${DOCDIR}
@dirrm ${HTMLDIR}/img
@dirrm ${HTMLDIR}
@dirrm ${PYSITELIB}/TMDA/pythonlib/email
@dirrm ${PYSITELIB}/TMDA/pythonlib
@dirrm ${PYSITELIB}/TMDA
