# mail/tmda/PLIST
@comment $NetBSD: PLIST,v 1.16 2009/06/14 18:04:45 joerg Exp $
bin/tmda-address
bin/tmda-check-address
bin/tmda-filter
bin/tmda-inject
bin/tmda-keygen
bin/tmda-ofmipd
bin/tmda-pending
bin/tmda-rfilter
bin/tmda-sendmail
${PYSITELIB}/TMDA/Address.py
${PYSITELIB}/TMDA/Address.pyc
${PYSITELIB}/TMDA/Auth.py
${PYSITELIB}/TMDA/Auth.pyc
${PYSITELIB}/TMDA/AutoResponse.py
${PYSITELIB}/TMDA/AutoResponse.pyc
${PYSITELIB}/TMDA/Cookie.py
${PYSITELIB}/TMDA/Cookie.pyc
${PYSITELIB}/TMDA/Defaults.py
${PYSITELIB}/TMDA/Defaults.pyc
${PYSITELIB}/TMDA/Deliver.py
${PYSITELIB}/TMDA/Deliver.pyc
${PYSITELIB}/TMDA/Errors.py
${PYSITELIB}/TMDA/Errors.pyc
${PYSITELIB}/TMDA/FilterParser.py
${PYSITELIB}/TMDA/FilterParser.pyc
${PYSITELIB}/TMDA/HMAC.py
${PYSITELIB}/TMDA/HMAC.pyc
${PYSITELIB}/TMDA/MTA.py
${PYSITELIB}/TMDA/MTA.pyc
${PYSITELIB}/TMDA/MessageLogger.py
${PYSITELIB}/TMDA/MessageLogger.pyc
${PYSITELIB}/TMDA/Pending.py
${PYSITELIB}/TMDA/Pending.pyc
${PYSITELIB}/TMDA/Queue/MaildirQueue.py
${PYSITELIB}/TMDA/Queue/MaildirQueue.pyc
${PYSITELIB}/TMDA/Queue/OriginalQueue.py
${PYSITELIB}/TMDA/Queue/OriginalQueue.pyc
${PYSITELIB}/TMDA/Queue/Queue.py
${PYSITELIB}/TMDA/Queue/Queue.pyc
${PYSITELIB}/TMDA/Queue/Util.py
${PYSITELIB}/TMDA/Queue/Util.pyc
${PYSITELIB}/TMDA/Queue/__init__.py
${PYSITELIB}/TMDA/Queue/__init__.pyc
${PYSITELIB}/TMDA/SMTP.py
${PYSITELIB}/TMDA/SMTP.pyc
${PYSITELIB}/TMDA/Util.py
${PYSITELIB}/TMDA/Util.pyc
${PYSITELIB}/TMDA/Version.py
${PYSITELIB}/TMDA/Version.pyc
${PYSITELIB}/TMDA/__init__.py
${PYSITELIB}/TMDA/__init__.pyc
${PYSITELIB}/TMDA/pythonlib/email/__init__.py
${PYSITELIB}/TMDA/pythonlib/email/__init__.pyc
${PYSITELIB}/TMDA/pythonlib/email/_parseaddr.py
${PYSITELIB}/TMDA/pythonlib/email/_parseaddr.pyc
${PYSITELIB}/TMDA/pythonlib/email/base64mime.py
${PYSITELIB}/TMDA/pythonlib/email/base64mime.pyc
${PYSITELIB}/TMDA/pythonlib/email/charset.py
${PYSITELIB}/TMDA/pythonlib/email/charset.pyc
${PYSITELIB}/TMDA/pythonlib/email/encoders.py
${PYSITELIB}/TMDA/pythonlib/email/encoders.pyc
${PYSITELIB}/TMDA/pythonlib/email/errors.py
${PYSITELIB}/TMDA/pythonlib/email/errors.pyc
${PYSITELIB}/TMDA/pythonlib/email/feedparser.py
${PYSITELIB}/TMDA/pythonlib/email/feedparser.pyc
${PYSITELIB}/TMDA/pythonlib/email/generator.py
${PYSITELIB}/TMDA/pythonlib/email/generator.pyc
${PYSITELIB}/TMDA/pythonlib/email/header.py
${PYSITELIB}/TMDA/pythonlib/email/header.pyc
${PYSITELIB}/TMDA/pythonlib/email/iterators.py
${PYSITELIB}/TMDA/pythonlib/email/iterators.pyc
${PYSITELIB}/TMDA/pythonlib/email/message.py
${PYSITELIB}/TMDA/pythonlib/email/message.pyc
${PYSITELIB}/TMDA/pythonlib/email/mime/__init__.py
${PYSITELIB}/TMDA/pythonlib/email/mime/__init__.pyc
${PYSITELIB}/TMDA/pythonlib/email/mime/application.py
${PYSITELIB}/TMDA/pythonlib/email/mime/application.pyc
${PYSITELIB}/TMDA/pythonlib/email/mime/audio.py
${PYSITELIB}/TMDA/pythonlib/email/mime/audio.pyc
${PYSITELIB}/TMDA/pythonlib/email/mime/base.py
${PYSITELIB}/TMDA/pythonlib/email/mime/base.pyc
${PYSITELIB}/TMDA/pythonlib/email/mime/image.py
${PYSITELIB}/TMDA/pythonlib/email/mime/image.pyc
${PYSITELIB}/TMDA/pythonlib/email/mime/message.py
${PYSITELIB}/TMDA/pythonlib/email/mime/message.pyc
${PYSITELIB}/TMDA/pythonlib/email/mime/multipart.py
${PYSITELIB}/TMDA/pythonlib/email/mime/multipart.pyc
${PYSITELIB}/TMDA/pythonlib/email/mime/nonmultipart.py
${PYSITELIB}/TMDA/pythonlib/email/mime/nonmultipart.pyc
${PYSITELIB}/TMDA/pythonlib/email/mime/text.py
${PYSITELIB}/TMDA/pythonlib/email/mime/text.pyc
${PYSITELIB}/TMDA/pythonlib/email/parser.py
${PYSITELIB}/TMDA/pythonlib/email/parser.pyc
${PYSITELIB}/TMDA/pythonlib/email/quoprimime.py
${PYSITELIB}/TMDA/pythonlib/email/quoprimime.pyc
${PYSITELIB}/TMDA/pythonlib/email/utils.py
${PYSITELIB}/TMDA/pythonlib/email/utils.pyc
${HTMLDIR}/AbandonedPages.html
${HTMLDIR}/AboutTmda.html
${HTMLDIR}/AutoAdminGroup.html
${HTMLDIR}/BadContent.html
${HTMLDIR}/CategoryCategory.html
${HTMLDIR}/CategoryHomepage.html
${HTMLDIR}/CategoryTemplate.html
${HTMLDIR}/ChallengeResponse.html
${HTMLDIR}/ClientConfiguration.html
${HTMLDIR}/ConfigurationOverview.html
${HTMLDIR}/ConfigurationVariables.html
${HTMLDIR}/DesktopEdition(2f)HowToConfigure.html
${HTMLDIR}/DesktopEdition.html
${HTMLDIR}/DrewRaines.html
${HTMLDIR}/EditedSystemPages.html
${HTMLDIR}/EventStats(2f)HitCounts.html
${HTMLDIR}/EventStats(2f)UserAgents.html
${HTMLDIR}/EventStats.html
${HTMLDIR}/FilterSources.html
${HTMLDIR}/FilterSpecification.html
${HTMLDIR}/FindPage.html
${HTMLDIR}/FortuneCookies.html
${HTMLDIR}/FrontPage.html
${HTMLDIR}/HeadHoncho.html
${HTMLDIR}/HelpContents.html
${HTMLDIR}/HelpForBeginners.html
${HTMLDIR}/HelpForDevelopers.html
${HTMLDIR}/HelpIndex.html
${HTMLDIR}/HelpMiscellaneous(2f)ExperimentalFeatures.html
${HTMLDIR}/HelpMiscellaneous(2f)FrequentlyAskedQuestions.html
${HTMLDIR}/HelpMiscellaneous.html
${HTMLDIR}/HelpOnAccessControlLists.html
${HTMLDIR}/HelpOnActions(2f)AttachFile.html
${HTMLDIR}/HelpOnActions.html
${HTMLDIR}/HelpOnAdministration.html
${HTMLDIR}/HelpOnAuthentication.html
${HTMLDIR}/HelpOnAutoAdmin.html
${HTMLDIR}/HelpOnCategories.html
${HTMLDIR}/HelpOnConfiguration(2f)ApacheVoodoo.html
${HTMLDIR}/HelpOnConfiguration(2f)CascadingStyleSheets.html
${HTMLDIR}/HelpOnConfiguration(2f)EmailSupport.html
${HTMLDIR}/HelpOnConfiguration(2f)FileAttachments.html
${HTMLDIR}/HelpOnConfiguration(2f)IntegratingWithApache.html
${HTMLDIR}/HelpOnConfiguration(2f)SecurityPolicy.html
${HTMLDIR}/HelpOnConfiguration(2f)SurgeProtection.html
${HTMLDIR}/HelpOnConfiguration(2f)UserPreferences.html
${HTMLDIR}/HelpOnConfiguration.html
${HTMLDIR}/HelpOnEditLocks.html
${HTMLDIR}/HelpOnEditing(2f)SubPages.html
${HTMLDIR}/HelpOnEditing.html
${HTMLDIR}/HelpOnFormatting.html
${HTMLDIR}/HelpOnGraphicalEditor.html
${HTMLDIR}/HelpOnHeadlines.html
${HTMLDIR}/HelpOnInstalling(2f)AolServer.html
${HTMLDIR}/HelpOnInstalling(2f)ApacheOnLinux.html
${HTMLDIR}/HelpOnInstalling(2f)ApacheOnLinuxFtp.html
${HTMLDIR}/HelpOnInstalling(2f)ApacheOnMacOsx.html
${HTMLDIR}/HelpOnInstalling(2f)ApacheOnWin32.html
${HTMLDIR}/HelpOnInstalling(2f)ApacheWithModPython.html
${HTMLDIR}/HelpOnInstalling(2f)BasicInstallation.html
${HTMLDIR}/HelpOnInstalling(2f)FastCgi.html
${HTMLDIR}/HelpOnInstalling(2f)InternetInformationServer.html
${HTMLDIR}/HelpOnInstalling(2f)MinimalWiki.html
${HTMLDIR}/HelpOnInstalling(2f)StandaloneServer.html
${HTMLDIR}/HelpOnInstalling(2f)TroubleShooting.html
${HTMLDIR}/HelpOnInstalling(2f)TwistedWeb.html
${HTMLDIR}/HelpOnInstalling(2f)WebLogic.html
${HTMLDIR}/HelpOnInstalling(2f)WikiInstanceCreation.html
${HTMLDIR}/HelpOnInstalling.html
${HTMLDIR}/HelpOnLanguages.html
${HTMLDIR}/HelpOnLinking.html
${HTMLDIR}/HelpOnLists.html
${HTMLDIR}/HelpOnMacros(2f)ImageLink.html
${HTMLDIR}/HelpOnMacros(2f)Include.html
${HTMLDIR}/HelpOnMacros(2f)MailTo.html
${HTMLDIR}/HelpOnMacros.html
${HTMLDIR}/HelpOnNavigation.html
${HTMLDIR}/HelpOnPackageInstaller.html
${HTMLDIR}/HelpOnPageCreation.html
${HTMLDIR}/HelpOnPageDeletion.html
${HTMLDIR}/HelpOnParsers(2f)ReStructuredText(2f)RstPrimer.html
${HTMLDIR}/HelpOnParsers(2f)ReStructuredText.html
${HTMLDIR}/HelpOnParsers.html
${HTMLDIR}/HelpOnProcessingInstructions.html
${HTMLDIR}/HelpOnProcessors.html
${HTMLDIR}/HelpOnRules.html
${HTMLDIR}/HelpOnSearching.html
${HTMLDIR}/HelpOnSkins.html
${HTMLDIR}/HelpOnSlideShows(2f)000_Introduction.html
${HTMLDIR}/HelpOnSlideShows(2f)100_Creating_the_slides.html
${HTMLDIR}/HelpOnSlideShows(2f)900_Last_but_not_least(3a)_Running_your_presentation(2e).html
${HTMLDIR}/HelpOnSlideShows.html
${HTMLDIR}/HelpOnSmileys.html
${HTMLDIR}/HelpOnSpellCheck.html
${HTMLDIR}/HelpOnSuperUser.html
${HTMLDIR}/HelpOnTables.html
${HTMLDIR}/HelpOnTemplates.html
${HTMLDIR}/HelpOnThemes.html
${HTMLDIR}/HelpOnUpdating.html
${HTMLDIR}/HelpOnUpdatingPython.html
${HTMLDIR}/HelpOnUserHandling.html
${HTMLDIR}/HelpOnUserPreferences.html
${HTMLDIR}/HelpOnVariables.html
${HTMLDIR}/HelpOnXmlPages.html
${HTMLDIR}/HelpTemplate.html
${HTMLDIR}/HomepageGroupsTemplate.html
${HTMLDIR}/HomepagePrivatePageTemplate.html
${HTMLDIR}/HomepageReadPageTemplate.html
${HTMLDIR}/HomepageReadWritePageTemplate.html
${HTMLDIR}/HomepageTemplate.html
${HTMLDIR}/InterWiki.html
${HTMLDIR}/JasonMastaler.html
${HTMLDIR}/LocalBadContent.html
${HTMLDIR}/LocalSpellingWords.html
${HTMLDIR}/MailingLists.html
${HTMLDIR}/MarkHorn.html
${HTMLDIR}/MissingHomePage.html
${HTMLDIR}/MissingPage.html
${HTMLDIR}/MoinMoin(2f)InstallDocs.html
${HTMLDIR}/MoinMoin(2f)TextFormatting.html
${HTMLDIR}/MoinMoin.html
${HTMLDIR}/MoinPagesEditorGroup.html
${HTMLDIR}/OrphanedPages.html
${HTMLDIR}/PageHits.html
${HTMLDIR}/PageSize.html
${HTMLDIR}/PreConfiguration.html
${HTMLDIR}/ProjectGroupsTemplate.html
${HTMLDIR}/ProjectTemplate.html
${HTMLDIR}/RandomPage.html
${HTMLDIR}/RecentChanges.html
${HTMLDIR}/ServerConfiguration.html
${HTMLDIR}/SiteNavigation.html
${HTMLDIR}/SlideShowHandOutTemplate.html
${HTMLDIR}/SlideShowTemplate.html
${HTMLDIR}/SlideTemplate.html
${HTMLDIR}/SpamAssassin.html
${HTMLDIR}/SquirrelMail.html
${HTMLDIR}/SyntaxReference.html
${HTMLDIR}/SystemInfo.html
${HTMLDIR}/SystemPagesGroup.html
${HTMLDIR}/SystemPagesInEnglishGroup.html
${HTMLDIR}/TemplateHowto.html
${HTMLDIR}/TitleIndex.html
${HTMLDIR}/TmdaBouncedDeletePendingHowto.html
${HTMLDIR}/TmdaCgi.html
${HTMLDIR}/TmdaCode.html
${HTMLDIR}/TmdaCommercialSupport.html
${HTMLDIR}/TmdaConfiguration.html
${HTMLDIR}/TmdaDocumentation.html
${HTMLDIR}/TmdaDonations.html
${HTMLDIR}/TmdaDownload.html
${HTMLDIR}/TmdaFaq.html
${HTMLDIR}/TmdaFeatures.html
${HTMLDIR}/TmdaFetchMailHowTo.html
${HTMLDIR}/TmdaFingerprint.html
${HTMLDIR}/TmdaHashCashHowto.html
${HTMLDIR}/TmdaHistory.html
${HTMLDIR}/TmdaHowtos.html
${HTMLDIR}/TmdaI18n.html
${HTMLDIR}/TmdaInUse.html
${HTMLDIR}/TmdaInstall.html
${HTMLDIR}/TmdaInstallation.html
${HTMLDIR}/TmdaLogos.html
${HTMLDIR}/TmdaOfmipdHowto.html
${HTMLDIR}/TmdaOneDotTwoDesign.html
${HTMLDIR}/TmdaOneDotTwoTodoList.html
${HTMLDIR}/TmdaPatchSubmissions.html
${HTMLDIR}/TmdaPendingAsMaildir.html
${HTMLDIR}/TmdaPendingAsMaildirHowto.html
${HTMLDIR}/TmdaPostfixPolicy.html
${HTMLDIR}/TmdaPressCoverage.html
${HTMLDIR}/TmdaProviders.html
${HTMLDIR}/TmdaRequirements.html
${HTMLDIR}/TmdaServices.html
${HTMLDIR}/TmdaSupport.html
${HTMLDIR}/TmdaTestimonials.html
${HTMLDIR}/TmdaTipsAndTricks.html
${HTMLDIR}/TmdaUpgrade.html
${HTMLDIR}/ToddJacobs.html
${HTMLDIR}/TroubleShooting.html
${HTMLDIR}/UserPreferences.html
${HTMLDIR}/VirtualDomainsHowto.html
${HTMLDIR}/WantedPages.html
${HTMLDIR}/WhyWikiWorks.html
${HTMLDIR}/WikiCourse(2f)BasicIntroduction(2f)000_What_is_a_Wiki(3f).html
${HTMLDIR}/WikiCourse(2f)BasicIntroduction(2f)050_First_Steps.html
${HTMLDIR}/WikiCourse(2f)BasicIntroduction(2f)060_How_do_I_navigate(3f).html
${HTMLDIR}/WikiCourse(2f)BasicIntroduction(2f)090_The_Editor.html
${HTMLDIR}/WikiCourse(2f)BasicIntroduction(2f)091_Tips_on_the_Editor.html
${HTMLDIR}/WikiCourse(2f)BasicIntroduction(2f)100_Text_markup.html
${HTMLDIR}/WikiCourse(2f)BasicIntroduction(2f)110_Headlines_and_Paragraphs.html
${HTMLDIR}/WikiCourse(2f)BasicIntroduction(2f)130_Lists.html
${HTMLDIR}/WikiCourse(2f)BasicIntroduction(2f)150_Tables.html
${HTMLDIR}/WikiCourse(2f)BasicIntroduction(2f)180_Macros.html
${HTMLDIR}/WikiCourse(2f)BasicIntroduction(2f)191_Creating_new_Pages.html
${HTMLDIR}/WikiCourse(2f)BasicIntroduction(2f)195_Actions.html
${HTMLDIR}/WikiCourse(2f)BasicIntroduction(2f)200_Wikiquette.html
${HTMLDIR}/WikiCourse(2f)BasicIntroduction(2f)210_Organisation_and_Structure.html
${HTMLDIR}/WikiCourse(2f)BasicIntroduction(2f)900_Exercises.html
${HTMLDIR}/WikiCourse(2f)BasicIntroduction.html
${HTMLDIR}/WikiCourse(2f)BasicIntroductionHandOut.html
${HTMLDIR}/WikiCourse.html
${HTMLDIR}/WikiHomePage.html
${HTMLDIR}/WikiName.html
${HTMLDIR}/WikiSandBox.html
${HTMLDIR}/WikiWikiWeb.html
${HTMLDIR}/WordIndex.html
${HTMLDIR}/XsltVersion.html
${HTMLDIR}/attachments/HelpOnActions(2f)AttachFile/unzip1.png
${HTMLDIR}/attachments/HelpOnActions(2f)AttachFile/unzip2.png
${HTMLDIR}/attachments/HelpOnActions(2f)AttachFile/unzip3a.png
${HTMLDIR}/attachments/HelpOnActions(2f)AttachFile/unzip3b.png
${HTMLDIR}/attachments/HelpOnActions(2f)AttachFile/unzip5.png
${HTMLDIR}/attachments/HelpOnActions(2f)AttachFile/unzip6.png
${HTMLDIR}/attachments/HelpOnActions(2f)AttachFile/unzip6b.png
${HTMLDIR}/attachments/HelpOnActions(2f)AttachFile/unzip7.png
${HTMLDIR}/attachments/HelpOnInstalling(2f)ApacheOnLinuxFtp/explore.py
${HTMLDIR}/attachments/HelpOnInstalling(2f)ApacheOnLinuxFtp/explore.pyc
${HTMLDIR}/attachments/HelpOnInstalling(2f)BasicInstallation/pythontest.cgi
${HTMLDIR}/attachments/HelpOnInstalling(2f)WikiInstanceCreation/createinstance.sh
${HTMLDIR}/attachments/HelpOnTables/MiniPage.png
${HTMLDIR}/attachments/TmdaFetchMailHowTo/xtmdaheader-1.1.c
${HTMLDIR}/attachments/TmdaFingerprint/tmda-fingerprint
${HTMLDIR}/attachments/TmdaPostfixPolicy/tmda-postfix-policy.sh
${HTMLDIR}/attachments/WikiSandBox/mytest.png
${HTMLDIR}/bounce_address(2e)py.html
${HTMLDIR}/deletepending_from_bouncemsg(2e)sh.html
${HTMLDIR}/delpenbounce.html
${HTMLDIR}/index.html
${HTMLDIR}/tmda(2d)fingerprint.html
${HTMLDIR}/tmda(2d)postfix(2d)policy(2e)sh.html
${DOCDIR}/CODENAMES
${DOCDIR}/COPYING
${DOCDIR}/CRYPTO
${DOCDIR}/ChangeLog
${DOCDIR}/INSTALL
${DOCDIR}/NEWS
${DOCDIR}/README
${DOCDIR}/THANKS
${DOCDIR}/UPGRADE
${EGDIR}/ChangeLog
${EGDIR}/collectaddys
${EGDIR}/def2html
${EGDIR}/printcdb
${EGDIR}/printdbm
${EGDIR}/sendit.sh
${EGDIR}/smtp-check-sender
${EGDIR}/tmda.el
${EGDIR}/tmda.spec
${EGDIR}/tofmipd.init
${EGDIR}/tofmipd.sysconfig
${EGDIR}/update-internaldomains
${EGDIR}/vadduser-tmda
${EGDIR}/vmailmgr-vdir.sh
${EGDIR}/vpopmail-vdir.sh
${EGDIR}/vtmdarc
${EGDIR}/wrapfd3.sh
${EGDIR}/dot-tmda/README
${EGDIR}/dot-tmda/config
${EGDIR}/dot-tmda/crypt_key
${EGDIR}/dot-tmda/filters/incoming
${EGDIR}/dot-tmda/filters/outgoing
${EGDIR}/dot-tmda/lists/blacklist
${EGDIR}/dot-tmda/lists/confirmed
${EGDIR}/dot-tmda/lists/whitelist
${EGDIR}/ofmipd-stunnel-daemons/README.txt
${EGDIR}/ofmipd-stunnel-daemons/stunnel.conf
${EGDIR}/ofmipd-stunnel-xinetd/README.txt
${EGDIR}/ofmipd-stunnel-xinetd/stunnel-wrapper
${EGDIR}/ofmipd-stunnel-xinetd/stunnel.conf
${EGDIR}/ofmipd-stunnel-xinetd/tmda-ofmipd-wrapper
${EGDIR}/ofmipd-stunnel-xinetd/xinetd.d-tmda-ofmipd-starttls
${SHAREDIR}/bounce_expired_dated.txt
${SHAREDIR}/bounce_fail_dated.txt
${SHAREDIR}/bounce_fail_keyword.txt
${SHAREDIR}/bounce_fail_sender.txt
${SHAREDIR}/bounce_incoming.txt
${SHAREDIR}/bounce_invalid_confirmation.txt
${SHAREDIR}/bounce_missing_pending.txt
${SHAREDIR}/confirm_accept.txt
${SHAREDIR}/confirm_request.txt
share/examples/rc.d/tofmipd
