# chat/ircu/PLIST
@comment $NetBSD: PLIST,v 1.4 2005/05/02 20:33:58 reed Exp $
man/man8/ircd.8
sbin/ircd
sbin/ircd-chkconf
sbin/ircd-mkpasswd
share/doc/ircd/README
share/doc/ircd/README.mkpasswd
share/doc/ircd/iauth.txt
share/doc/ircd/readme.asll
share/doc/ircd/readme.chroot
share/doc/ircd/readme.crules
share/doc/ircd/readme.cvs
share/doc/ircd/readme.features
share/doc/ircd/readme.gline
share/doc/ircd/readme.indent
share/doc/ircd/readme.jupe
share/doc/ircd/readme.log
share/doc/ircd/readme.who
share/doc/ircd/readme.www
share/examples/ircd/example.conf
share/examples/ircd/ircd.conf.sample
share/examples/rc.d/ircd
@dirrm share/examples/ircd
@dirrm share/doc/ircd
