# net/spreadlogd/PLIST
@comment $NetBSD: PLIST,v 1.1.1.1 2007/07/29 06:58:59 obache Exp $
sbin/spreadlogd
share/examples/rc.d/spreadlogd
share/examples/spreadlogd/spreadlogd.conf.sample
@dirrm share/examples/spreadlogd
