# databases/postgresql-server/pkg/PLIST
@comment $NetBSD: PLIST,v 1.1.1.1 2001/05/14 14:50:25 jlam Exp $
bin/postgres
bin/postmaster
etc/rc.d/pgsql
man/man1/postgres.1
man/man1/postmaster.1
share/postgresql/global.bki
share/postgresql/global.description
share/postgresql/pg_hba.conf.sample
share/postgresql/pg_ident.conf.sample
share/postgresql/postgresql.conf.sample
share/postgresql/template1.bki
share/postgresql/template1.description
@unexec rmdir -p %D/share/postgresql 2>/dev/null || true
