# sysutils/xentools20/PLIST
@comment $NetBSD: PLIST,v 1.5 2005/05/03 15:59:55 wiz Exp $
bin/miniterm
bin/xencons
bin/xenperf
include/xc.h
lib/libxc.so
lib/libxc.so.2.0
lib/libxc.so.2.0.0
lib/libxutil.a
lib/libxutil.so
lib/libxutil.so.2.0
lib/libxutil.so.2.0.0
sbin/netfix
sbin/xend
sbin/xenperf
sbin/xensv
sbin/xfrd
sbin/xm
${PYSITELIB}/xen/__init__.py
${PYSITELIB}/xen/lowlevel/__init__.py
${PYSITELIB}/xen/lowlevel/xc.so
${PYSITELIB}/xen/lowlevel/xu.so
${PYSITELIB}/xen/lowlevel/__init__.pyc
${PYSITELIB}/xen/util/__init__.py
${PYSITELIB}/xen/util/Brctl.py
${PYSITELIB}/xen/util/tempfile.py
${PYSITELIB}/xen/util/ip.py
${PYSITELIB}/xen/util/console_client.py
${PYSITELIB}/xen/util/__init__.pyc
${PYSITELIB}/xen/util/Brctl.pyc
${PYSITELIB}/xen/util/tempfile.pyc
${PYSITELIB}/xen/util/ip.pyc
${PYSITELIB}/xen/util/console_client.pyc
${PYSITELIB}/xen/xend/util.py
${PYSITELIB}/xen/xend/XendClient.py
${PYSITELIB}/xen/xend/XendProtocol.py
${PYSITELIB}/xen/xend/XendVnet.py
${PYSITELIB}/xen/xend/XendError.py
${PYSITELIB}/xen/xend/XendDomainInfo.py
${PYSITELIB}/xen/xend/XendDomain.py
${PYSITELIB}/xen/xend/XendLogging.py
${PYSITELIB}/xen/xend/XendRoot.py
${PYSITELIB}/xen/xend/XendDomainConfig.py
${PYSITELIB}/xen/xend/XendMigrate.py
${PYSITELIB}/xen/xend/PrettyPrint.py
${PYSITELIB}/xen/xend/XendConsole.py
${PYSITELIB}/xen/xend/sxp.py
${PYSITELIB}/xen/xend/Blkctl.py
${PYSITELIB}/xen/xend/XendNode.py
${PYSITELIB}/xen/xend/Vifctl.py
${PYSITELIB}/xen/xend/EventServer.py
${PYSITELIB}/xen/xend/EventTypes.py
${PYSITELIB}/xen/xend/XendDmesg.py
${PYSITELIB}/xen/xend/Args.py
${PYSITELIB}/xen/xend/__init__.py
${PYSITELIB}/xen/xend/XendAsynchProtocol.py
${PYSITELIB}/xen/xend/encode.py
${PYSITELIB}/xen/xend/XendDB.py
${PYSITELIB}/xen/xend/server/controller.py
${PYSITELIB}/xen/xend/server/SrvServer.py
${PYSITELIB}/xen/xend/server/SrvConsoleDir.py
${PYSITELIB}/xen/xend/server/SrvDeviceDir.py
${PYSITELIB}/xen/xend/server/__init__.py
${PYSITELIB}/xen/xend/server/messages.py
${PYSITELIB}/xen/xend/server/SrvBase.py
${PYSITELIB}/xen/xend/server/domain.py
${PYSITELIB}/xen/xend/server/SrvDir.py
${PYSITELIB}/xen/xend/server/netif.py
${PYSITELIB}/xen/xend/server/params.py
${PYSITELIB}/xen/xend/server/SrvDmesg.py
${PYSITELIB}/xen/xend/server/channel.py
${PYSITELIB}/xen/xend/server/console.py
${PYSITELIB}/xen/xend/server/SrvConsole.py
${PYSITELIB}/xen/xend/server/SrvRoot.py
${PYSITELIB}/xen/xend/server/SrvEventDir.py
${PYSITELIB}/xen/xend/server/SrvNode.py
${PYSITELIB}/xen/xend/server/SrvDaemon.py
${PYSITELIB}/xen/xend/server/SrvDomain.py
${PYSITELIB}/xen/xend/server/blkif.py
${PYSITELIB}/xen/xend/server/SrvDomainDir.py
${PYSITELIB}/xen/xend/server/SrvXendLog.py
${PYSITELIB}/xen/xend/server/SrvVnetDir.py
${PYSITELIB}/xen/xend/server/controller.pyc
${PYSITELIB}/xen/xend/server/SrvServer.pyc
${PYSITELIB}/xen/xend/server/SrvConsoleDir.pyc
${PYSITELIB}/xen/xend/server/SrvDeviceDir.pyc
${PYSITELIB}/xen/xend/server/__init__.pyc
${PYSITELIB}/xen/xend/server/messages.pyc
${PYSITELIB}/xen/xend/server/SrvBase.pyc
${PYSITELIB}/xen/xend/server/domain.pyc
${PYSITELIB}/xen/xend/server/SrvDir.pyc
${PYSITELIB}/xen/xend/server/netif.pyc
${PYSITELIB}/xen/xend/server/params.pyc
${PYSITELIB}/xen/xend/server/SrvDmesg.pyc
${PYSITELIB}/xen/xend/server/channel.pyc
${PYSITELIB}/xen/xend/server/console.pyc
${PYSITELIB}/xen/xend/server/SrvConsole.pyc
${PYSITELIB}/xen/xend/server/SrvRoot.pyc
${PYSITELIB}/xen/xend/server/SrvEventDir.pyc
${PYSITELIB}/xen/xend/server/SrvNode.pyc
${PYSITELIB}/xen/xend/server/SrvDaemon.pyc
${PYSITELIB}/xen/xend/server/SrvDomain.pyc
${PYSITELIB}/xen/xend/server/blkif.pyc
${PYSITELIB}/xen/xend/server/SrvDomainDir.pyc
${PYSITELIB}/xen/xend/server/SrvXendLog.pyc
${PYSITELIB}/xen/xend/server/SrvVnetDir.pyc
${PYSITELIB}/xen/xend/util.pyc
${PYSITELIB}/xen/xend/XendClient.pyc
${PYSITELIB}/xen/xend/XendProtocol.pyc
${PYSITELIB}/xen/xend/XendVnet.pyc
${PYSITELIB}/xen/xend/XendError.pyc
${PYSITELIB}/xen/xend/XendDomainInfo.pyc
${PYSITELIB}/xen/xend/XendDomain.pyc
${PYSITELIB}/xen/xend/XendLogging.pyc
${PYSITELIB}/xen/xend/XendRoot.pyc
${PYSITELIB}/xen/xend/XendDomainConfig.pyc
${PYSITELIB}/xen/xend/XendMigrate.pyc
${PYSITELIB}/xen/xend/PrettyPrint.pyc
${PYSITELIB}/xen/xend/XendConsole.pyc
${PYSITELIB}/xen/xend/sxp.pyc
${PYSITELIB}/xen/xend/Blkctl.pyc
${PYSITELIB}/xen/xend/XendNode.pyc
${PYSITELIB}/xen/xend/Vifctl.pyc
${PYSITELIB}/xen/xend/Args.pyc
${PYSITELIB}/xen/xend/EventServer.pyc
${PYSITELIB}/xen/xend/EventTypes.pyc
${PYSITELIB}/xen/xend/XendDmesg.pyc
${PYSITELIB}/xen/xend/__init__.pyc
${PYSITELIB}/xen/xend/XendAsynchProtocol.pyc
${PYSITELIB}/xen/xend/encode.pyc
${PYSITELIB}/xen/xend/XendDB.pyc
${PYSITELIB}/xen/sv/MigrateDomain.py
${PYSITELIB}/xen/sv/GenTabbed.py
${PYSITELIB}/xen/sv/HTMLBase.py
${PYSITELIB}/xen/sv/CreateDomain.py
${PYSITELIB}/xen/sv/Wizard.py
${PYSITELIB}/xen/sv/Main.py
${PYSITELIB}/xen/sv/TabView.py
${PYSITELIB}/xen/sv/DomInfo.py
${PYSITELIB}/xen/sv/util.py
${PYSITELIB}/xen/sv/SaveDomain.py
${PYSITELIB}/xen/sv/Daemon.py
${PYSITELIB}/xen/sv/NodeInfo.py
${PYSITELIB}/xen/sv/__init__.py
${PYSITELIB}/xen/sv/RestoreDomain.py
${PYSITELIB}/xen/sv/DomList.py
${PYSITELIB}/xen/sv/params.py
${PYSITELIB}/xen/sv/MigrateDomain.pyc
${PYSITELIB}/xen/sv/GenTabbed.pyc
${PYSITELIB}/xen/sv/HTMLBase.pyc
${PYSITELIB}/xen/sv/CreateDomain.pyc
${PYSITELIB}/xen/sv/Wizard.pyc
${PYSITELIB}/xen/sv/Main.pyc
${PYSITELIB}/xen/sv/TabView.pyc
${PYSITELIB}/xen/sv/DomInfo.pyc
${PYSITELIB}/xen/sv/util.pyc
${PYSITELIB}/xen/sv/SaveDomain.pyc
${PYSITELIB}/xen/sv/Daemon.pyc
${PYSITELIB}/xen/sv/NodeInfo.pyc
${PYSITELIB}/xen/sv/__init__.pyc
${PYSITELIB}/xen/sv/RestoreDomain.pyc
${PYSITELIB}/xen/sv/DomList.pyc
${PYSITELIB}/xen/sv/params.pyc
${PYSITELIB}/xen/xm/migrate.py
${PYSITELIB}/xen/xm/create.py
${PYSITELIB}/xen/xm/sysrq.py
${PYSITELIB}/xen/xm/help.py
${PYSITELIB}/xen/xm/opts.py
${PYSITELIB}/xen/xm/destroy.py
${PYSITELIB}/xen/xm/__init__.py
${PYSITELIB}/xen/xm/shutdown.py
${PYSITELIB}/xen/xm/main.py
${PYSITELIB}/xen/xm/migrate.pyc
${PYSITELIB}/xen/xm/create.pyc
${PYSITELIB}/xen/xm/sysrq.pyc
${PYSITELIB}/xen/xm/help.pyc
${PYSITELIB}/xen/xm/opts.pyc
${PYSITELIB}/xen/xm/destroy.pyc
${PYSITELIB}/xen/xm/__init__.pyc
${PYSITELIB}/xen/xm/shutdown.pyc
${PYSITELIB}/xen/xm/main.pyc
${PYSITELIB}/xen/__init__.pyc
share/examples/rc.d/xend
@dirrm ${PYSITELIB}/xen/lowlevel
@dirrm ${PYSITELIB}/xen/util
@dirrm ${PYSITELIB}/xen/xend/server
@dirrm ${PYSITELIB}/xen/xend
@dirrm ${PYSITELIB}/xen/sv
@dirrm ${PYSITELIB}/xen/xm
@dirrm ${PYSITELIB}/xen
