# net/tspc/PLIST
@comment $NetBSD: PLIST,v 1.6 2009/06/14 18:09:46 joerg Exp $
man/man5/tspc.conf.5
man/man8/tspc.8
sbin/tspc
share/examples/rc.d/tspc
share/examples/tspc/tspc.conf
share/tspc/template/README
share/tspc/template/checktunnel.sh
share/tspc/template/${TARGET}.sh
share/tspc/template/variables_environ
