# net/parpd/PLIST
@comment $NetBSD: PLIST,v 1.1.1.1 2008/10/09 00:24:28 joerg Exp $
man/man5/parpd.conf.5
man/man8/parpd.8
sbin/parpd
share/examples/rc.d/parpd
