# security/dirmngr/PLIST
@comment $NetBSD: PLIST,v 1.8 2006/04/17 07:07:33 jlam Exp $
bin/dirmngr
bin/dirmngr-client
info/dirmngr.info
libexec/dirmngr_ldap
sbin/runDirmngr
share/examples/rc.d/dirmngr
share/locale/de/LC_MESSAGES/dirmngr.mo
