# security/dirmngr/PLIST
@comment $NetBSD: PLIST,v 1.3 2005/01/04 13:40:38 shannonjr Exp $
bin/dirmngr
bin/dirmngr-client
libexec/dirmngr_ldap
sbin/runDirmngr
share/examples/rc.d/dirmngr
${PKGLOCALEDIR}/locale/de/LC_MESSAGES/dirmngr.mo
@exec ${MKDIR} %D/libdata/dirmngr/extra-certs
@dirrm libdata/dirmngr/extra-certs
@dirrm libdata/dirmngr
