# security/dirmngr/PLIST
@comment $NetBSD: PLIST,v 1.7 2006/03/31 23:56:29 jlam Exp $
bin/dirmngr
bin/dirmngr-client
info/dirmngr.info
libexec/dirmngr_ldap
sbin/runDirmngr
share/examples/rc.d/dirmngr
${PKGLOCALEDIR}/locale/de/LC_MESSAGES/dirmngr.mo
