# security/dirmngr/PLIST
@comment $NetBSD: PLIST,v 1.9 2006/10/11 11:09:48 shannonjr Exp $
bin/dirmngr
bin/dirmngr-client
info/dirmngr.info
libexec/dirmngr_ldap
man/man1/dirmngr-client.1
man/man1/dirmngr.1
sbin/runDirmngr
share/examples/rc.d/dirmngr
share/locale/de/LC_MESSAGES/dirmngr.mo
