# net/delegate/pkg/PLIST
@comment $NetBSD: PLIST,v 1.2 1999/02/23 18:28:35 rh Exp $
sbin/delegated
etc/rc.d/delegated.sh.sample
man/cat8/delegated.0
share/doc/delegate/tutorial-en.html
share/doc/delegate/tutorial-jp.html
@dirrm share/doc/delegate
