# sysutils/tits/PLIST
@comment $NetBSD: PLIST,v 1.3 2011/01/10 12:18:50 obache Exp $
man/man8/tits.8
man/cat8/tits.0
sbin/tits
share/examples/rc.d/tits
share/examples/tits.conf
